# tests/auto/OptionTooNew/2.34/good/Socket_LogLevelMax_emerg.target
[Unit]
Description=Foo
[Socket]
LogLevelMax=123

