# tests/auto/OptionTooNew/2.34/good/Unit_ConditionUser_abc.target
[Unit]
Description=Foo
Description=abc
ConditionUser=123

