# tests/auto/OptionTooNew/2.34/good/Unit_ConditionGroup_abc.target
[Unit]
Description=Foo
Description=abc
ConditionGroup=123

