# systemdlint/tests/auto/OptionTooNew/2.38/good/Socket_TemporaryFileSystem_abc.target
[Unit]
Description=Foo
[Socket]
TemporaryFileSystem=123

