# systemdlint/tests/auto/OptionTooNew/2.38/good/Swap_TemporaryFileSystem_abc.target
[Unit]
Description=Foo
[Swap]
TemporaryFileSystem=123

