# systemdlint/tests/auto/OptionTooNew/2.13/bad/Socket_RemoveOnStop_yes.target
[Unit]
Description=Foo
[Socket]
RemoveOnStop=123

