# systemdlint/tests/auto/OptionDeprecated/2.41/good/Mount_BlockIODeviceWeight_--dev--foo%201.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
BlockIODeviceWeight=123

