# systemdlint/tests/auto/OptionDeprecated/2.41/good/Mount_BlockIOWeight_1.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
BlockIOWeight=123

