# systemdlint/tests/auto/OptionDeprecated/2.41/good/Service_BlockIOAccounting_yes.target
[Unit]
Description=Foo
[Service]
BlockIOAccounting=123

