# systemdlint/tests/auto/OptionDeprecated/2.41/good/Service_BlockStartupIOWeight_1.target
[Unit]
Description=Foo
[Service]
BlockStartupIOWeight=123

