# systemdlint/tests/auto/OptionDeprecated/2.41/good/Service_BlockIODeviceWeight_--dev--foo%201.target
[Unit]
Description=Foo
[Service]
BlockIODeviceWeight=123

