# systemdlint/tests/auto/OptionDeprecated/2.41/good/Swap_CPUShares_2.target
[Unit]
Description=Foo
[Swap]
CPUShares=123

