# systemdlint/tests/auto/OptionDeprecated/2.41/good/Swap_BlockIOWriteBandwidthMax_--dev--foo%200K.target
[Unit]
Description=Foo
[Swap]
BlockIOWriteBandwidthMax=123

