# systemdlint/tests/auto/OptionDeprecated/2.41/good/Swap_BlockIOReadBandwidthMax_--dev--foo%200K.target
[Unit]
Description=Foo
[Swap]
BlockIOReadBandwidthMax=123

