# systemdlint/tests/auto/OptionDeprecated/2.41/good/Swap_StartupCPUShares_2.target
[Unit]
Description=Foo
[Swap]
StartupCPUShares=123

