# systemdlint/tests/auto/OptionTooNew/2.14/bad/Install_DefaultInstance_abc.target
[Unit]
Description=Foo
[Install]
DefaultInstance=123

