# systemdlint/tests/auto/OptionTooNew/2.14/bad/Service_RestartForceExitStatus_0.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=123

