# systemdlint/tests/auto/InvalidValue/2.49/bad/DHCPServer_RelayTarget_0--abc.target
[Unit]
Description=Foo
[DHCPServer]
RelayTarget=0/abc

