# systemdlint/tests/auto/InvalidValue/2.49/bad/DHCPServerStaticLease_MACAddress_00-11-22-GG-44-55.target
[Unit]
Description=Foo
[DHCPServerStaticLease]
MACAddress=00-11-22-GG-44-55

