# systemdlint/tests/auto/InvalidValue/2.49/bad/DHCPServer_RelayTarget_0.0.0..target
[Unit]
Description=Foo
[DHCPServer]
RelayTarget=0.0.0.

