# systemdlint/tests/auto/InvalidValue/2.49/bad/DHCPServerStaticLease_Address_0--abc.target
[Unit]
Description=Foo
[DHCPServerStaticLease]
Address=0/abc

