# systemdlint/tests/auto/InvalidValue/2.49/bad/IPv6AcceptRA_RouteMetric_4294967296.target
[Unit]
Description=Foo
[IPv6AcceptRA]
RouteMetric=4294967296

