# systemdlint/tests/auto/InvalidValue/2.49/bad/Upload_URL_0.0.0..target
[Unit]
Description=Foo
[Upload]
URL=0.0.0.

