# systemdlint/tests/auto/InvalidValue/2.49/bad/Upload_URL_255.255.255.255--33.target
[Unit]
Description=Foo
[Upload]
URL=255.255.255.255/33

