# systemdlint/tests/auto/InvalidValue/2.49/bad/Upload_URL_True.target
[Unit]
Description=Foo
[Upload]
URL=True

