# systemdlint/tests/auto/InvalidValue/2.49/bad/Upload_URL_abc.target
[Unit]
Description=Foo
[Upload]
URL=abc

