# systemdlint/tests/auto/InvalidValue/2.49/bad/Upload_URL_255.256.255.255.target
[Unit]
Description=Foo
[Upload]
URL=255.256.255.255

