# systemdlint/tests/auto/InvalidValue/2.49/bad/Upload_NetworkTimeoutSec_0.5d.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=0.5d

