# systemdlint/tests/auto/InvalidValue/2.49/bad/Upload_NetworkTimeoutSec_1p.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=1p

