# systemdlint/tests/auto/InvalidValue/2.49/bad/Upload_URL_1.1.1.1---1.target
[Unit]
Description=Foo
[Upload]
URL=1.1.1.1/-1

