# tests/auto/OptionTooNew/2.32/good/Manager_CtrlAltDelBurstAction_abc.target
[Unit]
Description=Foo
[Manager]
CtrlAltDelBurstAction=123

