# tests/auto/OptionTooNew/2.32/good/Service_CPUAccounting_yes.target
[Unit]
Description=Foo
[Service]
CPUAccounting=123

