# systemdlint/tests/auto/InvalidValue/2.47/good/BridgeMDB_MulticastGroupAddress_1.1.1.1--0.target
[Unit]
Description=Foo
[BridgeMDB]
MulticastGroupAddress=1.1.1.1/0

