# systemdlint/tests/auto/InvalidValue/2.47/good/BridgeMDB_MulticastGroupAddress_http%3A----a.b.com.target
[Unit]
Description=Foo
[BridgeMDB]
MulticastGroupAddress=http://a.b.com

