# systemdlint/tests/auto/InvalidValue/2.47/good/IPv6SendRA_DNS_abc.target
[Unit]
Description=Foo
[IPv6SendRA]
DNS=abc

