# systemdlint/tests/auto/InvalidValue/2.47/good/IPv6SendRA_RouterLifetimeSec_infinity.target
[Unit]
Description=Foo
[IPv6SendRA]
RouterLifetimeSec=infinity

