# systemdlint/tests/auto/InvalidValue/2.47/good/IPv6SendRA_RouterLifetimeSec_27y.target
[Unit]
Description=Foo
[IPv6SendRA]
RouterLifetimeSec=27y

