# systemdlint/tests/auto/InvalidValue/2.47/good/IPv6SendRA_OtherInformation_yes.target
[Unit]
Description=Foo
[IPv6SendRA]
OtherInformation=yes

