# systemdlint/tests/auto/InvalidValue/2.47/good/IPv6SendRA_DNSLifetimeSec_6h.target
[Unit]
Description=Foo
[IPv6SendRA]
DNSLifetimeSec=6h

