# systemdlint/tests/auto/InvalidValue/2.47/good/IPv6SendRA_DNSLifetimeSec_0days.target
[Unit]
Description=Foo
[IPv6SendRA]
DNSLifetimeSec=0days

