# systemdlint/tests/auto/InvalidValue/2.47/good/IPv6SendRA_RouterLifetimeSec_0days.target
[Unit]
Description=Foo
[IPv6SendRA]
RouterLifetimeSec=0days

