# systemdlint/tests/auto/InvalidValue/2.47/good/IPv6SendRA_Managed_yes.target
[Unit]
Description=Foo
[IPv6SendRA]
Managed=yes

