# systemdlint/tests/auto/InvalidValue/2.47/good/IPv6SendRA_EmitDomains_0.target
[Unit]
Description=Foo
[IPv6SendRA]
EmitDomains=0

