# systemdlint/tests/auto/InvalidValue/2.47/good/IPv6SendRA_DNSLifetimeSec_infinity.target
[Unit]
Description=Foo
[IPv6SendRA]
DNSLifetimeSec=infinity

