# systemdlint/tests/auto/InvalidValue/2.47/good/Network_DHCPv6PrefixDelegation_0.target
[Unit]
Description=Foo
[Network]
DHCPv6PrefixDelegation=0

