# systemdlint/tests/auto/InvalidValue/2.47/good/Network_IPv6SendRA_false.target
[Unit]
Description=Foo
[Network]
IPv6SendRA=false

