# systemdlint/tests/auto/InvalidValue/2.47/good/Network_IPv6SendRA_static.target
[Unit]
Description=Foo
[Network]
IPv6SendRA=static

