# systemdlint/tests/auto/InvalidValue/2.47/good/Service_ProcSubset_all.target
[Unit]
Description=Foo
[Service]
ProcSubset=all

