# systemdlint/tests/auto/InvalidValue/2.47/good/Service_MountImages_abc.target
[Unit]
Description=Foo
[Service]
MountImages=abc

