# systemdlint/tests/auto/InvalidValue/2.47/good/Service_ManagedOOMSwap_abc.target
[Unit]
Description=Foo
[Service]
ManagedOOMSwap=abc

