# systemdlint/tests/auto/InvalidValue/2.47/good/Service_SystemCallLog_abc.target
[Unit]
Description=Foo
[Service]
SystemCallLog=abc

