# systemdlint/tests/auto/InvalidValue/2.47/good/Service_ProtectProc_invisible.target
[Unit]
Description=Foo
[Service]
ProtectProc=invisible

