# systemdlint/tests/auto/InvalidValue/2.47/good/Socket_ManagedOOMMemoryPressure_abc.target
[Unit]
Description=Foo
[Socket]
ManagedOOMMemoryPressure=abc

