# systemdlint/tests/auto/InvalidValue/2.47/good/Socket_FlushPending_yes.target
[Unit]
Description=Foo
[Socket]
FlushPending=yes

