# systemdlint/tests/auto/InvalidValue/2.47/good/Socket_FlushPending_1.target
[Unit]
Description=Foo
[Socket]
FlushPending=1

