# systemdlint/tests/auto/InvalidValue/2.47/good/Socket_FlushPending_no.target
[Unit]
Description=Foo
[Socket]
FlushPending=no

