# systemdlint/tests/auto/InvalidValue/2.47/good/Socket_Timestamping_usec.target
[Unit]
Description=Foo
[Socket]
Timestamping=usec

