# systemdlint/tests/auto/InvalidValue/2.47/good/Socket_FlushPending_false.target
[Unit]
Description=Foo
[Socket]
FlushPending=false

