# systemdlint/tests/auto/InvalidValue/2.47/good/Socket_Timestamping_off.target
[Unit]
Description=Foo
[Socket]
Timestamping=off

