# systemdlint/tests/auto/InvalidValue/2.47/good/Swap_ProtectProc_noaccess.target
[Unit]
Description=Foo
[Swap]
ProtectProc=noaccess

