# systemdlint/tests/auto/InvalidValue/2.47/good/Swap_ManagedOOMSwap_abc.target
[Unit]
Description=Foo
[Swap]
ManagedOOMSwap=abc

