# systemdlint/tests/auto/InvalidValue/2.47/good/Swap_ProcSubset_pid.target
[Unit]
Description=Foo
[Swap]
ProcSubset=pid

