# systemdlint/tests/auto/InvalidValue/2.47/good/Swap_ManagedOOMMemoryPressureLimitPercent_abc.target
[Unit]
Description=Foo
[Swap]
ManagedOOMMemoryPressureLimitPercent=abc

