# systemdlint/tests/auto/InvalidValue/2.47/good/Swap_ProtectProc_ptraceable.target
[Unit]
Description=Foo
[Swap]
ProtectProc=ptraceable

