# systemdlint/tests/auto/InvalidValue/2.47/good/Swap_ProcSubset_all.target
[Unit]
Description=Foo
[Swap]
ProcSubset=all

