# systemdlint/tests/auto/InvalidValue/2.47/good/Timer_FixedRandomDelay_yes.target
[Unit]
Description=Foo
[Timer]
FixedRandomDelay=yes

