# systemdlint/tests/auto/InvalidValue/2.47/good/Timer_FixedRandomDelay_0.target
[Unit]
Description=Foo
[Timer]
FixedRandomDelay=0

