# systemdlint/tests/auto/InvalidValue/2.47/good/Timer_FixedRandomDelay_no.target
[Unit]
Description=Foo
[Timer]
FixedRandomDelay=no

