# systemdlint/tests/auto/InvalidValue/2.47/good/Timer_FixedRandomDelay_on.target
[Unit]
Description=Foo
[Timer]
FixedRandomDelay=on

