# systemdlint/tests/auto/InvalidValue/2.47/good/Timer_FixedRandomDelay_true.target
[Unit]
Description=Foo
[Timer]
FixedRandomDelay=true

