# systemdlint/tests/auto/InvalidValue/2.47/good/Timer_FixedRandomDelay_false.target
[Unit]
Description=Foo
[Timer]
FixedRandomDelay=false

