# systemdlint/tests/auto/InvalidValue/2.47/good/VXLAN_Independent_false.target
[Unit]
Description=Foo
[VXLAN]
Independent=false

