# systemdlint/tests/auto/InvalidValue/2.47/good/VXLAN_Independent_off.target
[Unit]
Description=Foo
[VXLAN]
Independent=off

