# systemdlint/tests/auto/InvalidValue/2.47/good/VXLAN_Independent_yes.target
[Unit]
Description=Foo
[VXLAN]
Independent=yes

