# systemdlint/tests/auto/InvalidValue/2.47/good/VXLAN_Independent_no.target
[Unit]
Description=Foo
[VXLAN]
Independent=no

