# systemdlint/tests/auto/OptionTooNew/2.49/good/Service_SocketBindAllow_abc.target
[Unit]
Description=Foo
[Service]
SocketBindAllow=123

