# systemdlint/tests/auto/OptionTooNew/2.49/good/Slice_SocketBindDeny_abc.target
[Unit]
Description=Foo
[Slice]
SocketBindDeny=123

