# systemdlint/tests/auto/OptionTooNew/2.49/good/Slice_SocketBindAllow_abc.target
[Unit]
Description=Foo
[Slice]
SocketBindAllow=123

