# systemdlint/tests/auto/OptionTooNew/2.49/good/Unit_PropagatesStopTo_syslog.service.target
[Unit]
Description=Foo
Description=abc
PropagatesStopTo=123

