# systemdlint/tests/auto/OptionTooNew/2.49/good/Upload_TrustedCertificateFile_abc.target
[Unit]
Description=Foo
[Upload]
TrustedCertificateFile=123

