# tests/auto/InvalidValue/2.43/good/Bridge_ProxyARPWiFi_no.target
[Unit]
Description=Foo
[Bridge]
ProxyARPWiFi=no

