# tests/auto/InvalidValue/2.43/good/Bridge_ProxyARPWiFi_true.target
[Unit]
Description=Foo
[Bridge]
ProxyARPWiFi=true

