# tests/auto/InvalidValue/2.43/good/Bridge_ProxyARP_off.target
[Unit]
Description=Foo
[Bridge]
ProxyARP=off

