# tests/auto/InvalidValue/2.43/good/Bridge_ProxyARP_1.target
[Unit]
Description=Foo
[Bridge]
ProxyARP=1

