# tests/auto/InvalidValue/2.43/good/Bridge_ProxyARPWiFi_yes.target
[Unit]
Description=Foo
[Bridge]
ProxyARPWiFi=yes

