# tests/auto/InvalidValue/2.43/good/Bridge_ProxyARPWiFi_off.target
[Unit]
Description=Foo
[Bridge]
ProxyARPWiFi=off

