# tests/auto/InvalidValue/2.43/good/Bridge_ProxyARP_0.target
[Unit]
Description=Foo
[Bridge]
ProxyARP=0

