# tests/auto/InvalidValue/2.43/good/DHCPv6_Anonymize_false.target
[Unit]
Description=Foo
[DHCPv6]
Anonymize=false

