# tests/auto/InvalidValue/2.43/good/DHCPv4_UseHostname_0.target
[Unit]
Description=Foo
[DHCPv4]
UseHostname=0

