# tests/auto/InvalidValue/2.43/good/DHCPv4_RapidCommit_on.target
[Unit]
Description=Foo
[DHCPv4]
RapidCommit=on

