# tests/auto/InvalidValue/2.43/good/DHCPv4_UseMTU_false.target
[Unit]
Description=Foo
[DHCPv4]
UseMTU=false

