# tests/auto/InvalidValue/2.43/good/DHCPv6_RapidCommit_true.target
[Unit]
Description=Foo
[DHCPv6]
RapidCommit=true

