# tests/auto/InvalidValue/2.43/good/DHCPv4_SendHostname_0.target
[Unit]
Description=Foo
[DHCPv4]
SendHostname=0

