# tests/auto/InvalidValue/2.43/good/DHCPv4_KeepConfiguration_abc.target
[Unit]
Description=Foo
[DHCPv4]
KeepConfiguration=abc

