# tests/auto/InvalidValue/2.43/good/DHCPv4_Anonymize_no.target
[Unit]
Description=Foo
[DHCPv4]
Anonymize=no

