# tests/auto/InvalidValue/2.43/good/DHCPv4_RequestBroadcast_abc.target
[Unit]
Description=Foo
[DHCPv4]
RequestBroadcast=abc

