# tests/auto/InvalidValue/2.43/good/DHCPv4_IAID_9999999999999.target
[Unit]
Description=Foo
[DHCPv4]
IAID=9999999999999

