# tests/auto/InvalidValue/2.43/good/DHCPv6_UseNTP_1.target
[Unit]
Description=Foo
[DHCPv6]
UseNTP=1

