# tests/auto/InvalidValue/2.43/good/DHCPv4_UseMTU_no.target
[Unit]
Description=Foo
[DHCPv4]
UseMTU=no

