# tests/auto/InvalidValue/2.43/good/DHCPv6_UseDNS_no.target
[Unit]
Description=Foo
[DHCPv6]
UseDNS=no

