# tests/auto/InvalidValue/2.43/good/DHCPv6_RapidCommit_0.target
[Unit]
Description=Foo
[DHCPv6]
RapidCommit=0

