# tests/auto/InvalidValue/2.43/good/DHCPv6_SendHostname_no.target
[Unit]
Description=Foo
[DHCPv6]
SendHostname=no

