# tests/auto/InvalidValue/2.43/good/DHCPv6_UseMTU_1.target
[Unit]
Description=Foo
[DHCPv6]
UseMTU=1

