# tests/auto/InvalidValue/2.43/good/DHCPv4_RouteMetric_9999999999999.target
[Unit]
Description=Foo
[DHCPv4]
RouteMetric=9999999999999

