# tests/auto/InvalidValue/2.43/good/DHCPv6_Anonymize_yes.target
[Unit]
Description=Foo
[DHCPv6]
Anonymize=yes

