# tests/auto/InvalidValue/2.43/good/DHCPv6_UseHostname_off.target
[Unit]
Description=Foo
[DHCPv6]
UseHostname=off

