# tests/auto/InvalidValue/2.43/good/FooOverUDP_Local_255.255.255.255.target
[Unit]
Description=Foo
[FooOverUDP]
Local=255.255.255.255

