# tests/auto/InvalidValue/2.43/good/FooOverUDP_PeerPort_65534.target
[Unit]
Description=Foo
[FooOverUDP]
PeerPort=65534

