# tests/auto/InvalidValue/2.43/good/GENEVE_GenericProtocolExtension_off.target
[Unit]
Description=Foo
[GENEVE]
GenericProtocolExtension=off

