# tests/auto/InvalidValue/2.43/good/GENEVE_IPDoNotFragment_1.target
[Unit]
Description=Foo
[GENEVE]
IPDoNotFragment=1

