# tests/auto/InvalidValue/2.43/good/GENEVE_IPDoNotFragment_true.target
[Unit]
Description=Foo
[GENEVE]
IPDoNotFragment=true

