# tests/auto/InvalidValue/2.43/good/GENEVE_GenericProtocolExtension_0.target
[Unit]
Description=Foo
[GENEVE]
GenericProtocolExtension=0

