# tests/auto/InvalidValue/2.43/good/GENEVE_GenericProtocolExtension_true.target
[Unit]
Description=Foo
[GENEVE]
GenericProtocolExtension=true

