# tests/auto/InvalidValue/2.43/good/GENEVE_IPDoNotFragment_yes.target
[Unit]
Description=Foo
[GENEVE]
IPDoNotFragment=yes

