# tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_25years.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=25years

