# tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_14msec.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=14msec

