# tests/auto/InvalidValue/2.43/good/Manager_NUMAPolicy_bind.target
[Unit]
Description=Foo
[Manager]
NUMAPolicy=bind

