# tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_2d.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=2d

