# tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_1day.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=1day

