# tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_6h.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=6h

