# tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_22weeks.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=22weeks

