# tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_7minutes.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=7minutes

