# tests/auto/InvalidValue/2.43/good/Manager_NUMAPolicy_local.target
[Unit]
Description=Foo
[Manager]
NUMAPolicy=local

