# tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_infinity.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=infinity

