# tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_23week.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=23week

