# tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_28%C2%B5s.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=28µs

