# tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_13months.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=13months

