# tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_10M.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=10M

