# tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_infinity.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=infinity

