# tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_12month.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=12month

