# tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_5hr.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=5hr

