# tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_3hours.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=3hours

