# tests/auto/InvalidValue/2.43/good/Manager_NUMAPolicy_interleave.target
[Unit]
Description=Foo
[Manager]
NUMAPolicy=interleave

