# tests/auto/InvalidValue/2.43/good/Match_Property_abc.target
[Unit]
Description=Foo
[Match]
Property=abc

