# tests/auto/InvalidValue/2.43/good/Network_DefaultRouteOnDevice_0.target
[Unit]
Description=Foo
[Network]
DefaultRouteOnDevice=0

