# tests/auto/InvalidValue/2.43/good/Network_DefaultRouteOnDevice_yes.target
[Unit]
Description=Foo
[Network]
DefaultRouteOnDevice=yes

