# tests/auto/InvalidValue/2.43/good/Network_DefaultRouteOnDevice_false.target
[Unit]
Description=Foo
[Network]
DefaultRouteOnDevice=false

