# tests/auto/InvalidValue/2.43/good/Route_TTLPropagate_1.target
[Unit]
Description=Foo
[Route]
TTLPropagate=1

