# tests/auto/InvalidValue/2.43/good/Route_FastOpenNoCookie_no.target
[Unit]
Description=Foo
[Route]
FastOpenNoCookie=no

