# tests/auto/InvalidValue/2.43/good/Route_TTLPropagate_yes.target
[Unit]
Description=Foo
[Route]
TTLPropagate=yes

