# tests/auto/InvalidValue/2.43/good/Route_TTLPropagate_false.target
[Unit]
Description=Foo
[Route]
TTLPropagate=false

