# tests/auto/InvalidValue/2.43/good/Route_TTLPropagate_on.target
[Unit]
Description=Foo
[Route]
TTLPropagate=on

