# tests/auto/InvalidValue/2.43/good/Route_TTLPropagate_0.target
[Unit]
Description=Foo
[Route]
TTLPropagate=0

