# tests/auto/InvalidValue/2.43/good/Route_TTLPropagate_true.target
[Unit]
Description=Foo
[Route]
TTLPropagate=true

