# tests/auto/InvalidValue/2.43/good/Route_FastOpenNoCookie_false.target
[Unit]
Description=Foo
[Route]
FastOpenNoCookie=false

