# tests/auto/InvalidValue/2.43/good/Route_FastOpenNoCookie_off.target
[Unit]
Description=Foo
[Route]
FastOpenNoCookie=off

