# tests/auto/InvalidValue/2.43/good/Route_FastOpenNoCookie_on.target
[Unit]
Description=Foo
[Route]
FastOpenNoCookie=on

