# tests/auto/InvalidValue/2.43/good/Service_DefaultMemoryLow_0B.target
[Unit]
Description=Foo
[Service]
DefaultMemoryLow=0B

