# tests/auto/InvalidValue/2.43/good/Service_NUMAPolicy_interleave.target
[Unit]
Description=Foo
[Service]
NUMAPolicy=interleave

