# tests/auto/InvalidValue/2.43/good/Service_NUMAPolicy_default.target
[Unit]
Description=Foo
[Service]
NUMAPolicy=default

