# tests/auto/InvalidValue/2.43/good/Service_TimeoutAbortSec_22weeks.target
[Unit]
Description=Foo
[Service]
TimeoutAbortSec=22weeks

