# tests/auto/InvalidValue/2.43/good/Service_TimeoutAbortSec_infinity.target
[Unit]
Description=Foo
[Service]
TimeoutAbortSec=infinity

