# tests/auto/InvalidValue/2.43/good/Service_DefaultMemoryMin_0B.target
[Unit]
Description=Foo
[Service]
DefaultMemoryMin=0B

