# tests/auto/InvalidValue/2.43/good/Socket_DefaultMemoryMin_9999999999999K.target
[Unit]
Description=Foo
[Socket]
DefaultMemoryMin=9999999999999K

