# tests/auto/InvalidValue/2.43/good/Socket_DefaultMemoryMin_0T.target
[Unit]
Description=Foo
[Socket]
DefaultMemoryMin=0T

