# tests/auto/InvalidValue/2.43/good/Socket_NUMAPolicy_preferred.target
[Unit]
Description=Foo
[Socket]
NUMAPolicy=preferred

