# tests/auto/InvalidValue/2.43/good/Socket_NUMAPolicy_bind.target
[Unit]
Description=Foo
[Socket]
NUMAPolicy=bind

