# tests/auto/InvalidValue/2.43/good/Swap_DefaultMemoryMin_9999999999999K.target
[Unit]
Description=Foo
[Swap]
DefaultMemoryMin=9999999999999K

