# tests/auto/InvalidValue/2.43/good/Swap_DefaultMemoryLow_0K.target
[Unit]
Description=Foo
[Swap]
DefaultMemoryLow=0K

