# tests/auto/InvalidValue/2.43/good/Swap_NUMAPolicy_preferred.target
[Unit]
Description=Foo
[Swap]
NUMAPolicy=preferred

