# tests/auto/InvalidValue/2.43/good/Swap_DefaultMemoryMin_0K.target
[Unit]
Description=Foo
[Swap]
DefaultMemoryMin=0K

