# tests/auto/InvalidValue/2.43/good/Swap_DefaultMemoryLow_9999999999999K.target
[Unit]
Description=Foo
[Swap]
DefaultMemoryLow=9999999999999K

