# tests/auto/InvalidValue/2.43/good/Swap_DefaultMemoryMin_0B.target
[Unit]
Description=Foo
[Swap]
DefaultMemoryMin=0B

