# tests/auto/InvalidValue/2.43/good/Swap_NUMAPolicy_default.target
[Unit]
Description=Foo
[Swap]
NUMAPolicy=default

