# tests/auto/InvalidValue/2.43/good/Tunnel_AssignToLoopback_no.target
[Unit]
Description=Foo
[Tunnel]
AssignToLoopback=no

