# tests/auto/InvalidValue/2.43/good/Tunnel_AssignToLoopback_0.target
[Unit]
Description=Foo
[Tunnel]
AssignToLoopback=0

