# tests/auto/InvalidValue/2.43/good/Tunnel_AssignToLoopback_yes.target
[Unit]
Description=Foo
[Tunnel]
AssignToLoopback=yes

