# tests/auto/InvalidValue/2.43/good/Tunnel_AssignToLoopback_off.target
[Unit]
Description=Foo
[Tunnel]
AssignToLoopback=off

