# tests/auto/InvalidValue/2.43/good/Tunnel_AssignToLoopback_on.target
[Unit]
Description=Foo
[Tunnel]
AssignToLoopback=on

