# tests/auto/InvalidValue/2.43/good/Tunnel_AssignToLoopback_1.target
[Unit]
Description=Foo
[Tunnel]
AssignToLoopback=1

