# tests/auto/InvalidValue/2.43/good/Unit_IPIngressFilterPath_abc.target
[Unit]
Description=Foo
Description=abc
IPIngressFilterPath=abc

