# tests/auto/InvalidValue/2.43/good/VXLAN_GenericProtocolExtension_yes.target
[Unit]
Description=Foo
[VXLAN]
GenericProtocolExtension=yes

