# tests/auto/InvalidValue/2.43/good/VXLAN_IPDoNotFragment_off.target
[Unit]
Description=Foo
[VXLAN]
IPDoNotFragment=off

