# tests/auto/InvalidValue/2.43/good/VXLAN_IPDoNotFragment_no.target
[Unit]
Description=Foo
[VXLAN]
IPDoNotFragment=no

