# tests/auto/InvalidValue/2.43/good/VXLAN_GenericProtocolExtension_no.target
[Unit]
Description=Foo
[VXLAN]
GenericProtocolExtension=no

