# tests/auto/InvalidValue/2.43/good/VXLAN_GenericProtocolExtension_1.target
[Unit]
Description=Foo
[VXLAN]
GenericProtocolExtension=1

