# tests/auto/InvalidValue/2.43/good/VXLAN_GenericProtocolExtension_off.target
[Unit]
Description=Foo
[VXLAN]
GenericProtocolExtension=off

