# systemdlint/tests/auto/InvalidValue/2.45/good/Address_AddPrefixRoute_on.target
[Unit]
Description=Foo
[Address]
Address=abc
AddPrefixRoute=on

