# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_TargetSec_27y.target
[Unit]
Description=Foo
[ControlledDelay]
TargetSec=27y

