# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_ECN_false.target
[Unit]
Description=Foo
[ControlledDelay]
ECN=false

