# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_17second.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=17second

