# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_TargetSec_8minute.target
[Unit]
Description=Foo
[ControlledDelay]
TargetSec=8minute

