# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_17second.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=17second

