# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_28%C2%B5s.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=28µs

