# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_26year.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=26year

