# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_1day.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=1day

