# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_1day.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=1day

