# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_12month.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=12month

