# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_Parent_root.target
[Unit]
Description=Foo
[ControlledDelay]
Parent=root

