# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_7minutes.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=7minutes

