# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_15ms.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=15ms

