# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_9min.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=9min

