# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_9min.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=9min

