# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_8minute.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=8minute

