# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_Handle_00.target
[Unit]
Description=Foo
[ControlledDelay]
Handle=00

