# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_19s.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=19s

