# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_25years.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=25years

