# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_14msec.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=14msec

