# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_21us.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=21us

