# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_TargetSec_9min.target
[Unit]
Description=Foo
[ControlledDelay]
TargetSec=9min

