# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_14msec.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=14msec

