# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_28%C2%B5s.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=28µs

