# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_TargetSec_25years.target
[Unit]
Description=Foo
[ControlledDelay]
TargetSec=25years

