# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_TargetSec_16seconds.target
[Unit]
Description=Foo
[ControlledDelay]
TargetSec=16seconds

