# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_27y.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=27y

