# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_22weeks.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=22weeks

