# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_ECN_0.target
[Unit]
Description=Foo
[ControlledDelay]
ECN=0

