# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_15ms.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=15ms

