# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_TargetSec_7minutes.target
[Unit]
Description=Foo
[ControlledDelay]
TargetSec=7minutes

