# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_21us.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=21us

