# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_ECN_true.target
[Unit]
Description=Foo
[ControlledDelay]
ECN=true

