# systemdlint/tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_4hour.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=4hour

