# systemdlint/tests/auto/InvalidValue/2.45/good/DHCPv4_SendDecline_1.target
[Unit]
Description=Foo
[DHCPv4]
SendDecline=1

