# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_ECN_false.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
ECN=false

