# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueing_CEThresholdSec_17second.target
[Unit]
Description=Foo
[FairQueueing]
CEThresholdSec=17second

