# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueing_CEThresholdSec_19s.target
[Unit]
Description=Foo
[FairQueueing]
CEThresholdSec=19s

