# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueing_Pacing_false.target
[Unit]
Description=Foo
[FairQueueing]
Pacing=false

