# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueing_CEThresholdSec_14msec.target
[Unit]
Description=Foo
[FairQueueing]
CEThresholdSec=14msec

