# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_20usec.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=20usec

