# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_Parent_root.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
Parent=root

