# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_6h.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=6h

