# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_9min.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=9min

