# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueing_Parent_clsact.target
[Unit]
Description=Foo
[FairQueueing]
Parent=clsact

