# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_infinity.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=infinity

