# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_2d.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=2d

