# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_Handle_00.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
Handle=00

