# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_27y.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=27y

