# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_7minutes.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=7minutes

