# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueing_MaximumRate_0K.target
[Unit]
Description=Foo
[FairQueueing]
MaximumRate=0K

