# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_5hr.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=5hr

