# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueing_CEThresholdSec_infinity.target
[Unit]
Description=Foo
[FairQueueing]
CEThresholdSec=infinity

