# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_15ms.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=15ms

