# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_ECN_yes.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
ECN=yes

