# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_Handle_fff.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
Handle=fff

