# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueing_FlowLimit_9999999999999.target
[Unit]
Description=Foo
[FairQueueing]
FlowLimit=9999999999999

