# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_13months.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=13months

