# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_28%C2%B5s.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=28µs

