# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_3hours.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=3hours

