# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_9min.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=9min

