# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueing_Pacing_off.target
[Unit]
Description=Foo
[FairQueueing]
Pacing=off

