# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_27y.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=27y

