# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueing_CEThresholdSec_28%C2%B5s.target
[Unit]
Description=Foo
[FairQueueing]
CEThresholdSec=28µs

