# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueing_InitialQuantum_9999999999999K.target
[Unit]
Description=Foo
[FairQueueing]
InitialQuantum=9999999999999K

