# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_6h.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=6h

