# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_ECN_0.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
ECN=0

