# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_10M.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=10M

