# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_25years.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=25years

