# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_ECN_on.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
ECN=on

