# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_6h.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=6h

