# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_8minute.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=8minute

