# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_12month.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=12month

