# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_2d.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=2d

