# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_11m.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=11m

