# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_MemoryLimit_0M.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
MemoryLimit=0M

