# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_25years.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=25years

