# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_17second.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=17second

