# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_1day.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=1day

