# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueing_FlowLimit_0.target
[Unit]
Description=Foo
[FairQueueing]
FlowLimit=0

