# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueing_Pacing_true.target
[Unit]
Description=Foo
[FairQueueing]
Pacing=true

