# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_23week.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=23week

