# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_19s.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=19s

