# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_26year.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=26year

