# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueing_CEThresholdSec_5hr.target
[Unit]
Description=Foo
[FairQueueing]
CEThresholdSec=5hr

