# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_7minutes.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=7minutes

