# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_4hour.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=4hour

