# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueing_OrphanMask_0.target
[Unit]
Description=Foo
[FairQueueing]
OrphanMask=0

