# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueing_Quantum_0G.target
[Unit]
Description=Foo
[FairQueueing]
Quantum=0G

