# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_5hr.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=5hr

