# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_24w.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=24w

