# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_0days.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=0days

