# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_9min.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=9min

