# systemdlint/tests/auto/InvalidValue/2.45/good/FairQueueing_Pacing_yes.target
[Unit]
Description=Foo
[FairQueueing]
Pacing=yes

