# systemdlint/tests/auto/InvalidValue/2.45/good/Link_AlternativeNamesPolicy_mac.target
[Unit]
Description=Foo
[Link]
MACAddress=abc
AlternativeNamesPolicy=mac

