# systemdlint/tests/auto/InvalidValue/2.45/good/Link_AlternativeNamesPolicy_slot.target
[Unit]
Description=Foo
[Link]
AlternativeNamesPolicy=slot

