# systemdlint/tests/auto/InvalidValue/2.45/good/Link_AlternativeName_abc.target
[Unit]
Description=Foo
[Link]
MACAddress=abc
AlternativeName=abc

