# systemdlint/tests/auto/InvalidValue/2.45/good/Link_AlternativeNamesPolicy_path.target
[Unit]
Description=Foo
[Link]
AlternativeNamesPolicy=path

