# systemdlint/tests/auto/InvalidValue/2.45/good/Match_PermanentMACAddress_abc.target
[Unit]
Description=Foo
[Match]
PermanentMACAddress=abc

