# systemdlint/tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_8minute.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=8minute

