# systemdlint/tests/auto/InvalidValue/2.45/good/NetworkEmulator_DuplicateRate_100%25.target
[Unit]
Description=Foo
[NetworkEmulator]
DuplicateRate=100%

