# systemdlint/tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_20usec.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=20usec

