# systemdlint/tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_11m.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=11m

