# systemdlint/tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_14msec.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=14msec

