# systemdlint/tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_9min.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=9min

