# systemdlint/tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_infinity.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=infinity

