# systemdlint/tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_12month.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=12month

