# systemdlint/tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_23week.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=23week

