# systemdlint/tests/auto/InvalidValue/2.45/good/NetworkEmulator_Parent_ingress.target
[Unit]
Description=Foo
[NetworkEmulator]
Parent=ingress

