# systemdlint/tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_6h.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=6h

