# systemdlint/tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_22weeks.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=22weeks

