# systemdlint/tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_26year.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=26year

