# systemdlint/tests/auto/InvalidValue/2.45/good/NetworkEmulator_Parent_clsact.target
[Unit]
Description=Foo
[NetworkEmulator]
Parent=clsact

