# systemdlint/tests/auto/InvalidValue/2.45/good/QDics_Parent_clsact.target
[Unit]
Description=Foo
[QDics]
Parent=clsact

