# systemdlint/tests/auto/InvalidValue/2.45/good/QDics_Parent_ingress.target
[Unit]
Description=Foo
[QDics]
Parent=ingress

