# systemdlint/tests/auto/InvalidValue/2.45/good/RoutingPolicyRule_User_abc.target
[Unit]
Description=Foo
[RoutingPolicyRule]
User=abc

