# systemdlint/tests/auto/InvalidValue/2.45/good/RoutingPolicyRule_SuppressPrefixLength_0.target
[Unit]
Description=Foo
[RoutingPolicyRule]
SuppressPrefixLength=0

