# systemdlint/tests/auto/InvalidValue/2.45/good/Service_ProtectClock_1.target
[Unit]
Description=Foo
[Service]
ProtectClock=1

