# systemdlint/tests/auto/InvalidValue/2.45/good/Service_ProtectClock_0.target
[Unit]
Description=Foo
[Service]
ProtectClock=0

