# systemdlint/tests/auto/InvalidValue/2.45/good/Service_ProtectClock_true.target
[Unit]
Description=Foo
[Service]
ProtectClock=true

