# systemdlint/tests/auto/InvalidValue/2.45/good/Service_ProtectClock_on.target
[Unit]
Description=Foo
[Service]
ProtectClock=on

