# systemdlint/tests/auto/InvalidValue/2.45/good/Socket_ProtectClock_true.target
[Unit]
Description=Foo
[Socket]
ProtectClock=true

