# systemdlint/tests/auto/InvalidValue/2.45/good/Socket_ProtectClock_off.target
[Unit]
Description=Foo
[Socket]
ProtectClock=off

