# systemdlint/tests/auto/InvalidValue/2.45/good/Socket_ProtectClock_0.target
[Unit]
Description=Foo
[Socket]
ProtectClock=0

