# systemdlint/tests/auto/InvalidValue/2.45/good/Socket_ProtectClock_on.target
[Unit]
Description=Foo
[Socket]
ProtectClock=on

