# systemdlint/tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_10M.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=10M

