# systemdlint/tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_11m.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=11m

