# systemdlint/tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_26year.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=26year

