# systemdlint/tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_16seconds.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=16seconds

