# systemdlint/tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_3hours.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=3hours

