# systemdlint/tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_9min.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=9min

