# systemdlint/tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_17second.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=17second

