# systemdlint/tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_1day.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=1day

