# systemdlint/tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_Parent_ingress.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
Parent=ingress

