# systemdlint/tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_22weeks.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=22weeks

