# systemdlint/tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_Handle_fff.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
Handle=fff

