# systemdlint/tests/auto/InvalidValue/2.45/good/Swap_ProtectClock_no.target
[Unit]
Description=Foo
[Swap]
ProtectClock=no

