# systemdlint/tests/auto/InvalidValue/2.45/good/Swap_ProtectClock_yes.target
[Unit]
Description=Foo
[Swap]
ProtectClock=yes

