# systemdlint/tests/auto/InvalidValue/2.45/good/Swap_ProtectClock_false.target
[Unit]
Description=Foo
[Swap]
ProtectClock=false

