# systemdlint/tests/auto/InvalidValue/2.45/good/Swap_ProtectClock_true.target
[Unit]
Description=Foo
[Swap]
ProtectClock=true

