# systemdlint/tests/auto/InvalidValue/2.45/good/Swap_ProtectClock_on.target
[Unit]
Description=Foo
[Swap]
ProtectClock=on

