# systemdlint/tests/auto/InvalidValue/2.45/good/TokenBucketFilter_LatencySec_13months.target
[Unit]
Description=Foo
[TokenBucketFilter]
LatencySec=13months

