# systemdlint/tests/auto/InvalidValue/2.45/good/TokenBucketFilter_MPUBytes_9999999999999K.target
[Unit]
Description=Foo
[TokenBucketFilter]
MPUBytes=9999999999999K

