# systemdlint/tests/auto/InvalidValue/2.45/good/TokenBucketFilter_PeakRate_0G.target
[Unit]
Description=Foo
[TokenBucketFilter]
PeakRate=0G

