# systemdlint/tests/auto/InvalidValue/2.45/good/TokenBucketFilter_PeakRate_9999999999999G.target
[Unit]
Description=Foo
[TokenBucketFilter]
PeakRate=9999999999999G

