# systemdlint/tests/auto/InvalidValue/2.45/good/TokenBucketFilter_LatencySec_4hour.target
[Unit]
Description=Foo
[TokenBucketFilter]
LatencySec=4hour

