# systemdlint/tests/auto/InvalidValue/2.45/good/TokenBucketFilter_Burst_0G.target
[Unit]
Description=Foo
[TokenBucketFilter]
Burst=0G

