# systemdlint/tests/auto/InvalidValue/2.45/good/TokenBucketFilter_MPUBytes_0M.target
[Unit]
Description=Foo
[TokenBucketFilter]
MPUBytes=0M

