# systemdlint/tests/auto/InvalidValue/2.45/good/TokenBucketFilter_LatencySec_11m.target
[Unit]
Description=Foo
[TokenBucketFilter]
LatencySec=11m

