# systemdlint/tests/auto/InvalidValue/2.45/good/TokenBucketFilter_LatencySec_28%C2%B5s.target
[Unit]
Description=Foo
[TokenBucketFilter]
LatencySec=28µs

