# systemdlint/tests/auto/InvalidValue/2.45/good/TokenBucketFilter_LatencySec_5hr.target
[Unit]
Description=Foo
[TokenBucketFilter]
LatencySec=5hr

