# systemdlint/tests/auto/InvalidValue/2.45/good/TokenBucketFilter_LimitSize_0M.target
[Unit]
Description=Foo
[TokenBucketFilter]
LimitSize=0M

