# systemdlint/tests/auto/InvalidValue/2.45/good/TokenBucketFilter_LatencySec_17second.target
[Unit]
Description=Foo
[TokenBucketFilter]
LatencySec=17second

