# systemdlint/tests/auto/InvalidValue/2.45/good/TokenBucketFilter_LatencySec_19s.target
[Unit]
Description=Foo
[TokenBucketFilter]
LatencySec=19s

