# systemdlint/tests/auto/InvalidValue/2.45/good/TokenBucketFilter_Parent_ingress.target
[Unit]
Description=Foo
[TokenBucketFilter]
Parent=ingress

