# systemdlint/tests/auto/InvalidValue/2.45/good/TokenBucketFilter_LatencySec_18sec.target
[Unit]
Description=Foo
[TokenBucketFilter]
LatencySec=18sec

