# systemdlint/tests/auto/InvalidValue/2.45/good/TokenBucketFilter_Handle_00.target
[Unit]
Description=Foo
[TokenBucketFilter]
Handle=00

