# systemdlint/tests/auto/InvalidValue/2.45/good/TrivialLinkEqualizer_Parent_ingress.target
[Unit]
Description=Foo
[TrivialLinkEqualizer]
Parent=ingress

