# systemdlint/tests/auto/InvalidValue/2.45/good/TrivialLinkEqualizer_Id_0.target
[Unit]
Description=Foo
[TrivialLinkEqualizer]
Id=0

