# systemdlint/tests/auto/InvalidValue/2.45/good/TrivialLinkEqualizer_Id_9999999999999.target
[Unit]
Description=Foo
[TrivialLinkEqualizer]
Id=9999999999999

