# systemdlint/tests/auto/InvalidValue/2.45/good/TrivialLinkEqualizer_Handle_fff.target
[Unit]
Description=Foo
[TrivialLinkEqualizer]
Handle=fff

