# tests/auto/OptionTooNew/2.47/bad/IPv6AcceptRA_PrefixAllowList_abc.target
[Unit]
Description=Foo
[IPv6AcceptRA]
PrefixAllowList=123

