# tests/auto/OptionTooNew/2.47/bad/IPv6AcceptRA_RouteDenyList_abc.target
[Unit]
Description=Foo
[IPv6AcceptRA]
RouteDenyList=123

