# tests/auto/OptionTooNew/2.47/bad/Link_TransmitQueueLength_0.target
[Unit]
Description=Foo
[Link]
TransmitQueueLength=123

