# tests/auto/OptionTooNew/2.47/bad/Manager_ManagerEnvironment_a%3Db.target
[Unit]
Description=Foo
[Manager]
ManagerEnvironment=123

