# tests/auto/OptionTooNew/2.47/bad/Network_RouteTable_abc.target
[Unit]
Description=Foo
[Network]
RouteTable=123

