# tests/auto/OptionTooNew/2.47/bad/OOM_DefaultMemoryPressureDurationSec_0.target
[Unit]
Description=Foo
[OOM]
DefaultMemoryPressureDurationSec=123

