# tests/auto/OptionTooNew/2.47/bad/Service_ManagedOOMPreference_none.target
[Unit]
Description=Foo
[Service]
ManagedOOMPreference=123

