# tests/auto/OptionTooNew/2.47/bad/Service_NoExecPaths_abc.target
[Unit]
Description=Foo
[Service]
NoExecPaths=123

