# tests/auto/OptionTooNew/2.47/bad/Slice_ManagedOOMPreference_none.target
[Unit]
Description=Foo
[Slice]
ManagedOOMPreference=123

