# tests/auto/OptionTooNew/2.47/bad/Swap_ManagedOOMPreference_none.target
[Unit]
Description=Foo
[Swap]
ManagedOOMPreference=123

