# tests/auto/OptionTooNew/2.47/bad/Unit_ConditionCPUFeature_abm.target
[Unit]
Description=Foo
Description=abc
ConditionCPUFeature=123

