# tests/auto/SettingRestricted/2.31/good/DHCP_UseRoutes_yesAnonymize%3Dtrue.target
[Unit]
Description=Foo
[DHCP]
UseRoutes=yes

