# tests/auto/SettingRestricted/2.31/good/DHCP_UseMTU_yesMTUBytes%3DNone-Anonymize%3Dtrue.target
[Unit]
Description=Foo
[DHCP]
UseMTU=yes

