# tests/auto/SettingRestricted/2.31/good/Network_PrimarySlave_yesMode%3Dbalance-alb.target
[Unit]
Description=Foo
[Network]
PrimarySlave=yes

