# tests/auto/OptionTooNew/2.10/bad/Manager_DefaultBlockIOAccounting_yes.target
[Unit]
Description=Foo
[Manager]
DefaultBlockIOAccounting=123

