# systemdlint/tests/auto/OptionDeprecated/2.42/bad/Mount_CPUShares_2.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
CPUShares=123

