# systemdlint/tests/auto/OptionDeprecated/2.42/bad/Service_CPUShares_2.target
[Unit]
Description=Foo
[Service]
CPUShares=123

