# systemdlint/tests/auto/OptionDeprecated/2.42/bad/Socket_BlockIOAccounting_yes.target
[Unit]
Description=Foo
[Socket]
BlockIOAccounting=123

