# systemdlint/tests/auto/OptionDeprecated/2.42/bad/Socket_BlockStartupIOWeight_1.target
[Unit]
Description=Foo
[Socket]
BlockStartupIOWeight=123

