# systemdlint/tests/auto/OptionDeprecated/2.42/bad/Socket_BlockIODeviceWeight_--dev--foo%201.target
[Unit]
Description=Foo
[Socket]
BlockIODeviceWeight=123

