# systemdlint/tests/auto/OptionDeprecated/2.42/bad/Socket_BlockIOWeight_1.target
[Unit]
Description=Foo
[Socket]
BlockIOWeight=123

