# systemdlint/tests/auto/InvalidValue/2.31/good/Address_HomeAddress_0.target
[Unit]
Description=Foo
[Address]
Address=abc
HomeAddress=0

