# systemdlint/tests/auto/InvalidValue/2.31/good/Address_DuplicateAddressDetection_false.target
[Unit]
Description=Foo
[Address]
Address=abc
DuplicateAddressDetection=false

