# systemdlint/tests/auto/InvalidValue/2.31/good/Address_Scope_link.target
[Unit]
Description=Foo
[Address]
Address=abc
Scope=link

