# systemdlint/tests/auto/InvalidValue/2.31/good/Address_ManageTemporaryAddress_false.target
[Unit]
Description=Foo
[Address]
Address=abc
ManageTemporaryAddress=false

