# systemdlint/tests/auto/InvalidValue/2.31/good/Address_PrefixRoute_off.target
[Unit]
Description=Foo
[Address]
Address=abc
PrefixRoute=off

