# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_28h.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=28h

