# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_22weeks.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=22weeks

