# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_27msec.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=27msec

