# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_DirectoryMode_007.target
[Unit]
Description=Foo
[Automount]
DirectoryMode=007

