# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_7seconds.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=7seconds

