# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_6days.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=6days

