# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_20usec.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=20usec

