# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_0days.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=0days

