# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_24years.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=24years

