# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_11m.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=11m

