# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_9month.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=9month

