# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_16week.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=16week

