# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_18sec.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=18sec

