# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_25s.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=25s

