# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_17minutes.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=17minutes

