# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_2d.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=2d

