# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_28%C2%B5s.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=28µs

