# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_10M.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=10M

