# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_11y.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=11y

