# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_10day.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=10day

