# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_14msec.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=14msec

