# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_12second.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=12second

