# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_21usec.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=21usec

