# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_5hour.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=5hour

