# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_8weeks.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=8weeks

