# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_15ms.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=15ms

