# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_26months.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=26months

