# systemdlint/tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_23week.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=23week

