# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_DynamicTransmitLoadBalancing_on.target
[Unit]
Description=Foo
[Bond]
DynamicTransmitLoadBalancing=on

