# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_LearnPacketIntervalSec_1.target
[Unit]
Description=Foo
[Bond]
LearnPacketIntervalSec=1

