# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_Mode_active-backup.target
[Unit]
Description=Foo
[Bond]
Mode=active-backup

