# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_AdActorSystemPriority_65535.target
[Unit]
Description=Foo
[Bond]
AdActorSystemPriority=65535

