# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_Mode_balance-xor.target
[Unit]
Description=Foo
[Bond]
Mode=balance-xor

