# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_ARPAllTargets_all.target
[Unit]
Description=Foo
[Bond]
ARPAllTargets=all

