# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_UpDelaySec_0.target
[Unit]
Description=Foo
[Bond]
UpDelaySec=0

