# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_UpDelaySec_9999999999999.target
[Unit]
Description=Foo
[Bond]
UpDelaySec=9999999999999

