# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_AdActorSystem_abc.target
[Unit]
Description=Foo
[Bond]
AdActorSystem=abc

