# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_DynamicTransmitLoadBalancing_false.target
[Unit]
Description=Foo
[Bond]
DynamicTransmitLoadBalancing=false

