# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_TransmitHashPolicy_encap3%2B4.target
[Unit]
Description=Foo
[Bond]
TransmitHashPolicy=encap3+4

