# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_AllSlavesActive_yes.target
[Unit]
Description=Foo
[Bond]
AllSlavesActive=yes

