# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_DynamicTransmitLoadBalancing_no.target
[Unit]
Description=Foo
[Bond]
DynamicTransmitLoadBalancing=no

