# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_FailOverMACPolicy_none.target
[Unit]
Description=Foo
[Bond]
FailOverMACPolicy=none

