# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_DynamicTransmitLoadBalancing_yes.target
[Unit]
Description=Foo
[Bond]
DynamicTransmitLoadBalancing=yes

