# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_DynamicTransmitLoadBalancing_off.target
[Unit]
Description=Foo
[Bond]
DynamicTransmitLoadBalancing=off

