# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_AdSelect_count.target
[Unit]
Description=Foo
[Bond]
AdSelect=count

