# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_ARPIPTargets_abc.target
[Unit]
Description=Foo
[Bond]
ARPIPTargets=abc

