# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_ResendIGMP_0.target
[Unit]
Description=Foo
[Bond]
ResendIGMP=0

