# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_Mode_802.3ad.target
[Unit]
Description=Foo
[Bond]
Mode=802.3ad

