# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_DownDelaySec_0.target
[Unit]
Description=Foo
[Bond]
DownDelaySec=0

