# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_PrimaryReselectPolicy_always.target
[Unit]
Description=Foo
[Bond]
PrimaryReselectPolicy=always

