# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_TransmitHashPolicy_layer3%2B4.target
[Unit]
Description=Foo
[Bond]
TransmitHashPolicy=layer3+4

