# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_ResendIGMP_255.target
[Unit]
Description=Foo
[Bond]
ResendIGMP=255

