# systemdlint/tests/auto/InvalidValue/2.31/good/Bond_MinLinks_0.target
[Unit]
Description=Foo
[Bond]
MinLinks=0

