# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_AllowPortToBeRoot_true.target
[Unit]
Description=Foo
[Bridge]
AllowPortToBeRoot=true

