# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_AllowPortToBeRoot_on.target
[Unit]
Description=Foo
[Bridge]
AllowPortToBeRoot=on

