# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_Learning_0.target
[Unit]
Description=Foo
[Bridge]
Learning=0

