# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_HairPin_1.target
[Unit]
Description=Foo
[Bridge]
HairPin=1

