# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_HairPin_yes.target
[Unit]
Description=Foo
[Bridge]
HairPin=yes

