# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_FastLeave_on.target
[Unit]
Description=Foo
[Bridge]
FastLeave=on

