# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_HairPin_no.target
[Unit]
Description=Foo
[Bridge]
HairPin=no

