# systemdlint/tests/auto/InvalidValue/2.31/good/BridgeFDB_MACAddress_00-11-22-33-44-55.target
[Unit]
Description=Foo
[BridgeFDB]
MACAddress=01:23:45:67:89:ab
MACAddress=00-11-22-33-44-55

