# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_Learning_yes.target
[Unit]
Description=Foo
[Bridge]
Learning=yes

