# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_UseBPDU_false.target
[Unit]
Description=Foo
[Bridge]
UseBPDU=false

