# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_NeighborSuppression_on.target
[Unit]
Description=Foo
[Bridge]
NeighborSuppression=on

