# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_AllowPortToBeRoot_false.target
[Unit]
Description=Foo
[Bridge]
AllowPortToBeRoot=false

