# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_NeighborSuppression_0.target
[Unit]
Description=Foo
[Bridge]
NeighborSuppression=0

