# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_MulticastFlood_yes.target
[Unit]
Description=Foo
[Bridge]
MulticastFlood=yes

