# systemdlint/tests/auto/InvalidValue/2.31/good/BridgeFDB_MACAddress_abc.target
[Unit]
Description=Foo
[BridgeFDB]
MACAddress=abc
MACAddress=abc

