# systemdlint/tests/auto/InvalidValue/2.31/good/BridgeFDB_VLANId_abc.target
[Unit]
Description=Foo
[BridgeFDB]
MACAddress=abc
VLANId=abc

