# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_MulticastToUnicast_yes.target
[Unit]
Description=Foo
[Bridge]
MulticastToUnicast=yes

