# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_FastLeave_yes.target
[Unit]
Description=Foo
[Bridge]
FastLeave=yes

