# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_UnicastFlood_no.target
[Unit]
Description=Foo
[Bridge]
UnicastFlood=no

