# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_FastLeave_0.target
[Unit]
Description=Foo
[Bridge]
FastLeave=0

