# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_HairPin_true.target
[Unit]
Description=Foo
[Bridge]
HairPin=true

