# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_MulticastToUnicast_on.target
[Unit]
Description=Foo
[Bridge]
MulticastToUnicast=on

