# systemdlint/tests/auto/InvalidValue/2.31/good/BridgeVLAN_EgressUntagged_abc.target
[Unit]
Description=Foo
[BridgeVLAN]
EgressUntagged=abc

