# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_MulticastToUnicast_no.target
[Unit]
Description=Foo
[Bridge]
MulticastToUnicast=no

