# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_NeighborSuppression_yes.target
[Unit]
Description=Foo
[Bridge]
NeighborSuppression=yes

