# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_MulticastFlood_off.target
[Unit]
Description=Foo
[Bridge]
MulticastFlood=off

