# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_AllowPortToBeRoot_off.target
[Unit]
Description=Foo
[Bridge]
AllowPortToBeRoot=off

