# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_MulticastFlood_on.target
[Unit]
Description=Foo
[Bridge]
MulticastFlood=on

