# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_FastLeave_true.target
[Unit]
Description=Foo
[Bridge]
FastLeave=true

