# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_MulticastFlood_true.target
[Unit]
Description=Foo
[Bridge]
MulticastFlood=true

