# systemdlint/tests/auto/InvalidValue/2.31/good/BridgeVLAN_VLAN_abc.target
[Unit]
Description=Foo
[BridgeVLAN]
VLAN=abc

