# systemdlint/tests/auto/InvalidValue/2.31/good/Bridge_MulticastFlood_0.target
[Unit]
Description=Foo
[Bridge]
MulticastFlood=0

