# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_27y.target
[Unit]
Description=Foo
[CAN]
RestartSec=27y

