# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_15min.target
[Unit]
Description=Foo
[CAN]
RestartSec=15min

