# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_10day.target
[Unit]
Description=Foo
[CAN]
RestartSec=10day

