# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_6h.target
[Unit]
Description=Foo
[CAN]
RestartSec=6h

