# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_14M.target
[Unit]
Description=Foo
[CAN]
RestartSec=14M

