# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_infinity.target
[Unit]
Description=Foo
[CAN]
RestartSec=infinity

