# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_Bitrate_8000.target
[Unit]
Description=Foo
[CAN]
Bitrate=8000

