# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_6days.target
[Unit]
Description=Foo
[CAN]
RestartSec=6days

