# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_8weeks.target
[Unit]
Description=Foo
[CAN]
RestartSec=8weeks

