# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_21us.target
[Unit]
Description=Foo
[CAN]
RestartSec=21us

