# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_1%C2%B5s.target
[Unit]
Description=Foo
[CAN]
RestartSec=1µs

