# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_Bitrate_0K.target
[Unit]
Description=Foo
[CAN]
Bitrate=0K

