# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_4ms.target
[Unit]
Description=Foo
[CAN]
RestartSec=4ms

