# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_17minutes.target
[Unit]
Description=Foo
[CAN]
RestartSec=17minutes

