# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_2us.target
[Unit]
Description=Foo
[CAN]
RestartSec=2us

