# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_7minutes.target
[Unit]
Description=Foo
[CAN]
RestartSec=7minutes

