# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_16week.target
[Unit]
Description=Foo
[CAN]
RestartSec=16week

