# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_25years.target
[Unit]
Description=Foo
[CAN]
RestartSec=25years

