# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_26months.target
[Unit]
Description=Foo
[CAN]
RestartSec=26months

