# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_9min.target
[Unit]
Description=Foo
[CAN]
RestartSec=9min

