# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_TripleSampling_off.target
[Unit]
Description=Foo
[CAN]
TripleSampling=off

