# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_5hr.target
[Unit]
Description=Foo
[CAN]
RestartSec=5hr

