# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_1day.target
[Unit]
Description=Foo
[CAN]
RestartSec=1day

