# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_0sec.target
[Unit]
Description=Foo
[CAN]
RestartSec=0sec

