# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_TripleSampling_0.target
[Unit]
Description=Foo
[CAN]
TripleSampling=0

