# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_28h.target
[Unit]
Description=Foo
[CAN]
RestartSec=28h

