# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_Bitrate_5000.target
[Unit]
Description=Foo
[CAN]
Bitrate=5000

