# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_13minute.target
[Unit]
Description=Foo
[CAN]
RestartSec=13minute

