# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_21usec.target
[Unit]
Description=Foo
[CAN]
RestartSec=21usec

