# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_19s.target
[Unit]
Description=Foo
[CAN]
RestartSec=19s

