# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_RestartSec_23m.target
[Unit]
Description=Foo
[CAN]
RestartSec=23m

