# systemdlint/tests/auto/InvalidValue/2.31/good/CAN_TripleSampling_no.target
[Unit]
Description=Foo
[CAN]
TripleSampling=no

