# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_28h.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=28h

