# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_SendHostname_true.target
[Unit]
Description=Foo
[DHCP]
SendHostname=true

