# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_EmitDNS_no.target
[Unit]
Description=Foo
[DHCPServer]
EmitDNS=no

