# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_11m.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=11m

