# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_EmitRouter_true.target
[Unit]
Description=Foo
[DHCPServer]
EmitRouter=true

