# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_Hostname_abc.target
[Unit]
Description=Foo
[DHCP]
Hostname=abc

