# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_UseHostname_no.target
[Unit]
Description=Foo
[DHCP]
UseHostname=no

