# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_27y.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=27y

