# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_28h.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=28h

