# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_17second.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=17second

