# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_16seconds.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=16seconds

