# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_UseRoutes_on.target
[Unit]
Description=Foo
[DHCP]
UseRoutes=on

