# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_Anonymize_0.target
[Unit]
Description=Foo
[DHCP]
Anonymize=0

