# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_UseTimezone_true.target
[Unit]
Description=Foo
[DHCP]
UseTimezone=true

