# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_9month.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=9month

