# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_UseNTP_on.target
[Unit]
Description=Foo
[DHCP]
UseNTP=on

