# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_25s.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=25s

