# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_UseNTP_0.target
[Unit]
Description=Foo
[DHCP]
UseNTP=0

