# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_ListenPort_1.target
[Unit]
Description=Foo
[DHCP]
ListenPort=1

