# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_EmitDNS_on.target
[Unit]
Description=Foo
[DHCPServer]
EmitDNS=on

