# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_EmitDNS_false.target
[Unit]
Description=Foo
[DHCPServer]
EmitDNS=false

