# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_UseMTU_off.target
[Unit]
Description=Foo
[DHCP]
UseMTU=off

