# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_1%C2%B5s.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=1µs

