# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_CriticalConnection_yes.target
[Unit]
Description=Foo
[DHCP]
CriticalConnection=yes

