# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_2us.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=2us

