# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_1%C2%B5s.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=1µs

