# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_RapidCommit_on.target
[Unit]
Description=Foo
[DHCP]
RapidCommit=on

