# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_24w.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=24w

