# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_12second.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=12second

