# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_20w.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=20w

