# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_10day.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=10day

