# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_26months.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=26months

