# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_CriticalConnection_no.target
[Unit]
Description=Foo
[DHCP]
CriticalConnection=no

