# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_23m.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=23m

