# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_UseDomains_false.target
[Unit]
Description=Foo
[DHCP]
UseDomains=false

