# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_4ms.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=4ms

