# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_UseTimezone_yes.target
[Unit]
Description=Foo
[DHCP]
UseTimezone=yes

