# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_13months.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=13months

