# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_RapidCommit_true.target
[Unit]
Description=Foo
[DHCP]
RapidCommit=true

