# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_Anonymize_no.target
[Unit]
Description=Foo
[DHCP]
Anonymize=no

