# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_24years.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=24years

