# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_12month.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=12month

