# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_5hour.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=5hour

