# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_ClientIdentifier_duid.target
[Unit]
Description=Foo
[DHCP]
ClientIdentifier=duid

