# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_15ms.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=15ms

