# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_CriticalConnection_on.target
[Unit]
Description=Foo
[DHCP]
CriticalConnection=on

