# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_15ms.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=15ms

