# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_EmitNTP_off.target
[Unit]
Description=Foo
[DHCPServer]
EmitNTP=off

