# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_0days.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=0days

