# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_3hours.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=3hours

