# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_UseDomains_yes.target
[Unit]
Description=Foo
[DHCP]
UseDomains=yes

