# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_15min.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=15min

