# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_27msec.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=27msec

