# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_23week.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=23week

