# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_28%C2%B5s.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=28µs

