# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_UseHostname_true.target
[Unit]
Description=Foo
[DHCP]
UseHostname=true

