# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_UseDNS_yes.target
[Unit]
Description=Foo
[DHCP]
UseDNS=yes

