# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_6days.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=6days

