# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_RouteTable_4294967295.target
[Unit]
Description=Foo
[DHCP]
RouteTable=4294967295

