# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_EmitNTP_true.target
[Unit]
Description=Foo
[DHCPServer]
EmitNTP=true

