# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_IAID_9999999999999.target
[Unit]
Description=Foo
[DHCP]
IAID=9999999999999

