# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_UseTimezone_no.target
[Unit]
Description=Foo
[DHCP]
UseTimezone=no

