# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_9min.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=9min

