# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_13minute.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=13minute

