# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_EmitNTP_on.target
[Unit]
Description=Foo
[DHCPServer]
EmitNTP=on

