# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_25years.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=25years

