# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_Anonymize_true.target
[Unit]
Description=Foo
[DHCP]
Anonymize=true

