# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_8weeks.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=8weeks

