# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_22weeks.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=22weeks

