# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_25s.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=25s

