# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_27msec.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=27msec

