# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_26months.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=26months

