# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_RouteMetric_9999999999999.target
[Unit]
Description=Foo
[DHCP]
RouteMetric=9999999999999

