# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_UseDNS_no.target
[Unit]
Description=Foo
[DHCP]
UseDNS=no

