# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_14msec.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=14msec

