# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_7seconds.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=7seconds

