# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_21usec.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=21usec

