# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_ForceDHCPv6PDOtherInformation_yes.target
[Unit]
Description=Foo
[DHCP]
ForceDHCPv6PDOtherInformation=yes

