# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_UseTimezone_false.target
[Unit]
Description=Foo
[DHCP]
UseTimezone=false

