# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_EmitTimezone_on.target
[Unit]
Description=Foo
[DHCPServer]
EmitTimezone=on

