# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_18year.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=18year

