# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_UseMTU_false.target
[Unit]
Description=Foo
[DHCP]
UseMTU=false

