# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_20w.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=20w

