# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_UseRoutes_0.target
[Unit]
Description=Foo
[DHCP]
UseRoutes=0

