# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_UseRoutes_false.target
[Unit]
Description=Foo
[DHCP]
UseRoutes=false

