# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_5hour.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=5hour

