# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_28%C2%B5s.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=28µs

