# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_19d.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=19d

