# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_23m.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=23m

