# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_Anonymize_off.target
[Unit]
Description=Foo
[DHCP]
Anonymize=off

