# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_UseTimezone_0.target
[Unit]
Description=Foo
[DHCP]
UseTimezone=0

