# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_RapidCommit_false.target
[Unit]
Description=Foo
[DHCP]
RapidCommit=false

