# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_9month.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=9month

