# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DNS_abc.target
[Unit]
Description=Foo
[DHCPServer]
DNS=abc

