# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_0sec.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=0sec

