# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_10day.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=10day

