# systemdlint/tests/auto/InvalidValue/2.31/good/DHCP_CriticalConnection_false.target
[Unit]
Description=Foo
[DHCP]
CriticalConnection=false

