# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_21us.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=21us

