# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_27y.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=27y

