# systemdlint/tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_15min.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=15min

