# systemdlint/tests/auto/InvalidValue/2.31/good/FooOverUDP_Protocol_gre.target
[Unit]
Description=Foo
[FooOverUDP]
Protocol=gre

