# systemdlint/tests/auto/InvalidValue/2.31/good/FooOverUDP_Encapsulation_abc.target
[Unit]
Description=Foo
[FooOverUDP]
Encapsulation=abc

