# systemdlint/tests/auto/InvalidValue/2.31/good/GENEVE_Remote_1.1.1.1.target
[Unit]
Description=Foo
[GENEVE]
Remote=1.1.1.1

