# systemdlint/tests/auto/InvalidValue/2.31/good/GENEVE_DestinationPort_1.target
[Unit]
Description=Foo
[GENEVE]
DestinationPort=1

