# systemdlint/tests/auto/InvalidValue/2.31/good/GENEVE_UDPChecksum_0.target
[Unit]
Description=Foo
[GENEVE]
UDPChecksum=0

