# systemdlint/tests/auto/InvalidValue/2.31/good/GENEVE_UDPChecksum_off.target
[Unit]
Description=Foo
[GENEVE]
UDPChecksum=off

