# systemdlint/tests/auto/InvalidValue/2.31/good/GENEVE_UDP6ZeroChecksumRx_1.target
[Unit]
Description=Foo
[GENEVE]
UDP6ZeroChecksumRx=1

