# systemdlint/tests/auto/InvalidValue/2.31/good/GENEVE_DestinationPort_65534.target
[Unit]
Description=Foo
[GENEVE]
DestinationPort=65534

