# systemdlint/tests/auto/InvalidValue/2.31/good/Install_Alias_abc.target
[Unit]
Description=Foo
[Install]
Alias=abc

