# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6AddressLabel_Label_abc.target
[Unit]
Description=Foo
[IPv6AddressLabel]
Label=abc
Prefix=abc
Label=abc

