# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6Prefix_AddressAutoconfiguration_1.target
[Unit]
Description=Foo
[IPv6Prefix]
AddressAutoconfiguration=1

