# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_RouterLifetimeSec_5hour.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
RouterLifetimeSec=5hour

