# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6AcceptRA_UseDNS_1.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseDNS=1

