# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_OtherInformation_0.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
OtherInformation=0

