# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_EmitDomains_1.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
EmitDomains=1

