# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6AcceptRA_UseAutonomousPrefix_off.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseAutonomousPrefix=off

