# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6AcceptRA_UseOnLinkPrefix_off.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseOnLinkPrefix=off

