# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_DNSLifetimeSec_11y.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
DNSLifetimeSec=11y

