# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6AcceptRA_UseDomains_true.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseDomains=true

