# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6Prefix_ValidLifetimeSec_20w.target
[Unit]
Description=Foo
[IPv6Prefix]
ValidLifetimeSec=20w

