# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_DNS_abc.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
DNS=abc

