# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6Prefix_PreferredLifetimeSec_16week.target
[Unit]
Description=Foo
[IPv6Prefix]
PreferredLifetimeSec=16week

