# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6AcceptRA_UseDNS_0.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseDNS=0

