# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6Prefix_PreferredLifetimeSec_22hr.target
[Unit]
Description=Foo
[IPv6Prefix]
PreferredLifetimeSec=22hr

