# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6AcceptRA_UseDomains_false.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseDomains=false

