# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_RouterLifetimeSec_8weeks.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
RouterLifetimeSec=8weeks

