# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6Prefix_PreferredLifetimeSec_28h.target
[Unit]
Description=Foo
[IPv6Prefix]
PreferredLifetimeSec=28h

