# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_RouterLifetimeSec_14msec.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
RouterLifetimeSec=14msec

