# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_EmitDNS_off.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
EmitDNS=off

