# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6Prefix_ValidLifetimeSec_27msec.target
[Unit]
Description=Foo
[IPv6Prefix]
ValidLifetimeSec=27msec

