# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_RouterLifetimeSec_19d.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
RouterLifetimeSec=19d

