# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_OtherInformation_yes.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
OtherInformation=yes

