# systemdlint/tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_DNSLifetimeSec_17minutes.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
DNSLifetimeSec=17minutes

