# systemdlint/tests/auto/InvalidValue/2.31/good/IPVLAN_Mode_L2.target
[Unit]
Description=Foo
[IPVLAN]
Mode=L2

