# systemdlint/tests/auto/InvalidValue/2.31/good/IPVLAN_Flags_bridge.target
[Unit]
Description=Foo
[IPVLAN]
Flags=bridge

