# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxLevelSyslog_info.target
[Unit]
Description=Foo
[Journal]
MaxLevelSyslog=info

