# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_RuntimeMaxFiles_9999999999999E.target
[Unit]
Description=Foo
[Journal]
RuntimeMaxFiles=9999999999999E

