# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxRetentionSec_0h.target
[Unit]
Description=Foo
[Journal]
MaxRetentionSec=0h

