# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_14msec.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=14msec

