# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxLevelWall_err.target
[Unit]
Description=Foo
[Journal]
MaxLevelWall=err

