# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SplitMode_uid.target
[Unit]
Description=Foo
[Journal]
SplitMode=uid

