# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_ForwardToKMsg_on.target
[Unit]
Description=Foo
[Journal]
ForwardToKMsg=on

