# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_ForwardToSyslog_yes.target
[Unit]
Description=Foo
[Journal]
ForwardToSyslog=yes

