# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_20usec.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=20usec

