# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxFileSec_0year.target
[Unit]
Description=Foo
[Journal]
MaxFileSec=0year

