# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_ForwardToSyslog_off.target
[Unit]
Description=Foo
[Journal]
ForwardToSyslog=off

