# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_10M.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=10M

