# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_24w.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=24w

