# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxLevelWall_debug.target
[Unit]
Description=Foo
[Journal]
MaxLevelWall=debug

