# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_19s.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=19s

