# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_LineMax_79M.target
[Unit]
Description=Foo
[Journal]
LineMax=79M

