# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_LineMax_9999999999999G.target
[Unit]
Description=Foo
[Journal]
LineMax=9999999999999G

