# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_Compress_no.target
[Unit]
Description=Foo
[Journal]
Compress=no

