# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_ForwardToSyslog_0.target
[Unit]
Description=Foo
[Journal]
ForwardToSyslog=0

