# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SystemMaxFiles_9999999999999P.target
[Unit]
Description=Foo
[Journal]
SystemMaxFiles=9999999999999P

