# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SystemMaxUse_0M.target
[Unit]
Description=Foo
[Journal]
SystemMaxUse=0M

