# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_ReadKMsg_false.target
[Unit]
Description=Foo
[Journal]
ReadKMsg=false

