# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_ForwardToWall_false.target
[Unit]
Description=Foo
[Journal]
ForwardToWall=false

