# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SystemMaxFiles_0G.target
[Unit]
Description=Foo
[Journal]
SystemMaxFiles=0G

