# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_ForwardToWall_no.target
[Unit]
Description=Foo
[Journal]
ForwardToWall=no

