# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_ReadKMsg_off.target
[Unit]
Description=Foo
[Journal]
ReadKMsg=off

