# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_Storage_auto.target
[Unit]
Description=Foo
[Journal]
Storage=auto

