# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_ReadKMsg_0.target
[Unit]
Description=Foo
[Journal]
ReadKMsg=0

