# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_7minutes.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=7minutes

