# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_infinity.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=infinity

