# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxLevelWall_notice.target
[Unit]
Description=Foo
[Journal]
MaxLevelWall=notice

