# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxLevelWall_warning.target
[Unit]
Description=Foo
[Journal]
MaxLevelWall=warning

