# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_ForwardToConsole_0.target
[Unit]
Description=Foo
[Journal]
ForwardToConsole=0

