# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_16seconds.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=16seconds

