# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_RateLimitBurst_9999999999999.target
[Unit]
Description=Foo
[Journal]
RateLimitBurst=9999999999999

