# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_17second.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=17second

