# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxLevelStore_0.target
[Unit]
Description=Foo
[Journal]
MaxLevelStore=0

