# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SystemKeepFree_9999999999999G.target
[Unit]
Description=Foo
[Journal]
SystemKeepFree=9999999999999G

