# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_12month.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=12month

