# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_11m.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=11m

