# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_ForwardToKMsg_1.target
[Unit]
Description=Foo
[Journal]
ForwardToKMsg=1

