# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxLevelKMsg_info.target
[Unit]
Description=Foo
[Journal]
MaxLevelKMsg=info

