# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_1day.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=1day

