# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_ForwardToKMsg_false.target
[Unit]
Description=Foo
[Journal]
ForwardToKMsg=false

