# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_26year.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=26year

