# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxLevelWall_crit.target
[Unit]
Description=Foo
[Journal]
MaxLevelWall=crit

