# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_ForwardToKMsg_no.target
[Unit]
Description=Foo
[Journal]
ForwardToKMsg=no

