# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_Storage_volatile.target
[Unit]
Description=Foo
[Journal]
Storage=volatile

