# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_RateLimitIntervalSec_9999999999999min.target
[Unit]
Description=Foo
[Journal]
RateLimitIntervalSec=9999999999999min

