# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxLevelConsole_warning.target
[Unit]
Description=Foo
[Journal]
MaxLevelConsole=warning

