# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxRetentionSec_9999999999999m.target
[Unit]
Description=Foo
[Journal]
MaxRetentionSec=9999999999999m

