# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_ReadKMsg_no.target
[Unit]
Description=Foo
[Journal]
ReadKMsg=no

