# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxLevelKMsg_err.target
[Unit]
Description=Foo
[Journal]
MaxLevelKMsg=err

