# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_Compress_on.target
[Unit]
Description=Foo
[Journal]
Compress=on

