# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_ForwardToWall_1.target
[Unit]
Description=Foo
[Journal]
ForwardToWall=1

