# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxLevelConsole_emerg.target
[Unit]
Description=Foo
[Journal]
MaxLevelConsole=emerg

