# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxRetentionSec_0month.target
[Unit]
Description=Foo
[Journal]
MaxRetentionSec=0month

