# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_15ms.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=15ms

