# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_ForwardToWall_true.target
[Unit]
Description=Foo
[Journal]
ForwardToWall=true

