# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxLevelKMsg_7.target
[Unit]
Description=Foo
[Journal]
MaxLevelKMsg=7

