# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxLevelKMsg_emerg.target
[Unit]
Description=Foo
[Journal]
MaxLevelKMsg=emerg

