# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxLevelConsole_7.target
[Unit]
Description=Foo
[Journal]
MaxLevelConsole=7

