# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_MaxLevelWall_7.target
[Unit]
Description=Foo
[Journal]
MaxLevelWall=7

