# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_Seal_0.target
[Unit]
Description=Foo
[Journal]
Seal=0

