# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_ReadKMsg_on.target
[Unit]
Description=Foo
[Journal]
ReadKMsg=on

