# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_RuntimeKeepFree_0E.target
[Unit]
Description=Foo
[Journal]
RuntimeKeepFree=0E

