# systemdlint/tests/auto/InvalidValue/2.31/good/Journal_Seal_false.target
[Unit]
Description=Foo
[Journal]
Seal=false

