# systemdlint/tests/auto/InvalidValue/2.31/good/L2TP_PeerTunnelId_4294967295.target
[Unit]
Description=Foo
[L2TP]
PeerTunnelId=4294967295

