# systemdlint/tests/auto/InvalidValue/2.31/good/L2TP_UDPSourcePort_1.target
[Unit]
Description=Foo
[L2TP]
UDPSourcePort=1

