# systemdlint/tests/auto/InvalidValue/2.31/good/Link_MACAddressPolicy_abc.target
[Unit]
Description=Foo
[Link]
MACAddress=abc
MACAddressPolicy=abc

