# systemdlint/tests/auto/InvalidValue/2.31/good/Link_AllMulticast_false.target
[Unit]
Description=Foo
[Link]
AllMulticast=false

