# systemdlint/tests/auto/InvalidValue/2.31/good/Link_Unmanaged_no.target
[Unit]
Description=Foo
[Link]
MACAddress=01:23:45:67:89:ab 00-11-22-33-44-55 AABB.CCDD.EEFF
Unmanaged=no

