# systemdlint/tests/auto/InvalidValue/2.31/good/Link_AllMulticast_on.target
[Unit]
Description=Foo
[Link]
AllMulticast=on

