# systemdlint/tests/auto/InvalidValue/2.31/good/Link_ARP_false.target
[Unit]
Description=Foo
[Link]
ARP=false

