# systemdlint/tests/auto/InvalidValue/2.31/good/Link_RequiredForOnline_true.target
[Unit]
Description=Foo
[Link]
RequiredForOnline=true

