# systemdlint/tests/auto/InvalidValue/2.31/good/Link_Unmanaged_false.target
[Unit]
Description=Foo
[Link]
Unmanaged=false

