# systemdlint/tests/auto/InvalidValue/2.31/good/Link_MTUBytes_9999999999999M.target
[Unit]
Description=Foo
[Link]
MACAddress=abc
MTUBytes=9999999999999M

