# systemdlint/tests/auto/InvalidValue/2.31/good/Link_Unmanaged_1.target
[Unit]
Description=Foo
[Link]
Unmanaged=1

