# systemdlint/tests/auto/InvalidValue/2.31/good/Link_MTUBytes_9999999999999K.target
[Unit]
Description=Foo
[Link]
MTUBytes=9999999999999K

