# systemdlint/tests/auto/InvalidValue/2.31/good/Link_Multicast_true.target
[Unit]
Description=Foo
[Link]
Multicast=true

