# systemdlint/tests/auto/InvalidValue/2.31/good/Link_Multicast_no.target
[Unit]
Description=Foo
[Link]
Multicast=no

