# systemdlint/tests/auto/InvalidValue/2.31/good/Link_RequiredForOnline_off.target
[Unit]
Description=Foo
[Link]
RequiredForOnline=off

