# systemdlint/tests/auto/InvalidValue/2.31/good/Link_Unmanaged_no.target
[Unit]
Description=Foo
[Link]
MACAddress=abc
Unmanaged=no

