# systemdlint/tests/auto/InvalidValue/2.31/good/Link_Multicast_yes.target
[Unit]
Description=Foo
[Link]
Multicast=yes

