# systemdlint/tests/auto/InvalidValue/2.31/good/Link_Unmanaged_0.target
[Unit]
Description=Foo
[Link]
Unmanaged=0

