# systemdlint/tests/auto/InvalidValue/2.31/good/Link_AllMulticast_0.target
[Unit]
Description=Foo
[Link]
AllMulticast=0

