# systemdlint/tests/auto/InvalidValue/2.31/good/Link_Multicast_false.target
[Unit]
Description=Foo
[Link]
Multicast=false

