# systemdlint/tests/auto/InvalidValue/2.31/good/Link_MACAddress_abc.target
[Unit]
Description=Foo
[Link]
MACAddress=abc
MACAddress=abc

