# systemdlint/tests/auto/InvalidValue/2.31/good/Login_SuspendKeyIgnoreInhibited_yes.target
[Unit]
Description=Foo
[Login]
SuspendKeyIgnoreInhibited=yes

