# systemdlint/tests/auto/InvalidValue/2.31/good/Login_PowerKeyIgnoreInhibited_off.target
[Unit]
Description=Foo
[Login]
PowerKeyIgnoreInhibited=off

