# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_12second.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=12second

