# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitchDocked_halt.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchDocked=halt

