# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_12second.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=12second

