# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_15ms.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=15ms

