# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_14M.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=14M

