# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_23week.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=23week

