# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_14M.target
[Unit]
Description=Foo
[Login]
IdleActionSec=14M

