# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_5hr.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=5hr

