# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_20w.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=20w

