# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitch_halt.target
[Unit]
Description=Foo
[Login]
HandleLidSwitch=halt

