# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_9month.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=9month

