# systemdlint/tests/auto/InvalidValue/2.31/good/Login_PowerKeyIgnoreInhibited_on.target
[Unit]
Description=Foo
[Login]
PowerKeyIgnoreInhibited=on

