# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_21usec.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=21usec

