# systemdlint/tests/auto/InvalidValue/2.31/good/Login_KillUserProcesses_yes.target
[Unit]
Description=Foo
[Login]
KillUserProcesses=yes

