# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleAction_halt.target
[Unit]
Description=Foo
[Login]
IdleAction=halt

