# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_20usec.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=20usec

