# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleAction_poweroff.target
[Unit]
Description=Foo
[Login]
IdleAction=poweroff

