# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_28h.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=28h

