# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitchDocked_hybrid-sleep.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchDocked=hybrid-sleep

