# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_25years.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=25years

