# systemdlint/tests/auto/InvalidValue/2.31/good/Login_PowerKeyIgnoreInhibited_false.target
[Unit]
Description=Foo
[Login]
PowerKeyIgnoreInhibited=false

