# systemdlint/tests/auto/InvalidValue/2.31/good/Login_SessionsMax_0M.target
[Unit]
Description=Foo
[Login]
SessionsMax=0M

