# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandlePowerKey_kexec.target
[Unit]
Description=Foo
[Login]
HandlePowerKey=kexec

