# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_15ms.target
[Unit]
Description=Foo
[Login]
IdleActionSec=15ms

