# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitchDocked_suspend-then-hibernate.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchDocked=suspend-then-hibernate

