# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleAction_suspend.target
[Unit]
Description=Foo
[Login]
IdleAction=suspend

