# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_14M.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=14M

