# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitch_suspend.target
[Unit]
Description=Foo
[Login]
HandleLidSwitch=suspend

