# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_7seconds.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=7seconds

