# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleHibernateKey_lock.target
[Unit]
Description=Foo
[Login]
HandleHibernateKey=lock

