# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_4hour.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=4hour

