# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitchExternalPower_hibernate.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchExternalPower=hibernate

