# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_19s.target
[Unit]
Description=Foo
[Login]
IdleActionSec=19s

