# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_15ms.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=15ms

