# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleHibernateKey_ignore.target
[Unit]
Description=Foo
[Login]
HandleHibernateKey=ignore

