# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_26months.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=26months

