# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_20w.target
[Unit]
Description=Foo
[Login]
IdleActionSec=20w

