# systemdlint/tests/auto/InvalidValue/2.31/good/Login_SuspendKeyIgnoreInhibited_off.target
[Unit]
Description=Foo
[Login]
SuspendKeyIgnoreInhibited=off

