# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitch_kexec.target
[Unit]
Description=Foo
[Login]
HandleLidSwitch=kexec

