# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_24years.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=24years

