# systemdlint/tests/auto/InvalidValue/2.31/good/Login_SuspendKeyIgnoreInhibited_1.target
[Unit]
Description=Foo
[Login]
SuspendKeyIgnoreInhibited=1

