# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_infinity.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=infinity

