# systemdlint/tests/auto/InvalidValue/2.31/good/Login_SessionsMax_0T.target
[Unit]
Description=Foo
[Login]
SessionsMax=0T

