# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_20w.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=20w

