# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_0sec.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=0sec

