# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleAction_lock.target
[Unit]
Description=Foo
[Login]
IdleAction=lock

