# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitchExternalPower_ignore.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchExternalPower=ignore

