# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_11y.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=11y

