# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_26months.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=26months

