# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitchDocked_reboot.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchDocked=reboot

