# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_1day.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=1day

