# systemdlint/tests/auto/InvalidValue/2.31/good/Login_RuntimeDirectorySize_4096.target
[Unit]
Description=Foo
[Login]
RuntimeDirectorySize=4096

