# systemdlint/tests/auto/InvalidValue/2.31/good/Login_PowerKeyIgnoreInhibited_true.target
[Unit]
Description=Foo
[Login]
PowerKeyIgnoreInhibited=true

