# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitchDocked_kexec.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchDocked=kexec

