# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_8weeks.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=8weeks

