# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_6days.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=6days

