# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_19d.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=19d

