# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_16seconds.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=16seconds

