# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_23week.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=23week

