# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_27msec.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=27msec

