# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_28h.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=28h

