# systemdlint/tests/auto/InvalidValue/2.31/good/Login_PowerKeyIgnoreInhibited_1.target
[Unit]
Description=Foo
[Login]
PowerKeyIgnoreInhibited=1

