# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_17minutes.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=17minutes

