# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_5hour.target
[Unit]
Description=Foo
[Login]
IdleActionSec=5hour

