# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleAction_kexec.target
[Unit]
Description=Foo
[Login]
IdleAction=kexec

