# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_24w.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=24w

