# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_13months.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=13months

