# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleHibernateKey_hibernate.target
[Unit]
Description=Foo
[Login]
HandleHibernateKey=hibernate

