# systemdlint/tests/auto/InvalidValue/2.31/good/Login_PowerKeyIgnoreInhibited_no.target
[Unit]
Description=Foo
[Login]
PowerKeyIgnoreInhibited=no

