# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_20w.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=20w

