# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitorsMax_0M.target
[Unit]
Description=Foo
[Login]
InhibitorsMax=0M

