# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitorsMax_0G.target
[Unit]
Description=Foo
[Login]
InhibitorsMax=0G

