# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_12month.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=12month

