# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleSuspendKey_lock.target
[Unit]
Description=Foo
[Login]
HandleSuspendKey=lock

