# systemdlint/tests/auto/InvalidValue/2.31/good/Login_KillExcludeUsers_abc.target
[Unit]
Description=Foo
[Login]
KillExcludeUsers=abc

