# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_21usec.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=21usec

