# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_1%C2%B5s.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=1µs

