# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_22hr.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=22hr

