# systemdlint/tests/auto/InvalidValue/2.31/good/Login_RemoveIPC_1.target
[Unit]
Description=Foo
[Login]
RemoveIPC=1

