# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_10M.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=10M

