# systemdlint/tests/auto/InvalidValue/2.31/good/Login_RemoveIPC_no.target
[Unit]
Description=Foo
[Login]
RemoveIPC=no

