# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_17minutes.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=17minutes

