# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_10day.target
[Unit]
Description=Foo
[Login]
IdleActionSec=10day

