# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_27y.target
[Unit]
Description=Foo
[Login]
IdleActionSec=27y

