# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitchDocked_poweroff.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchDocked=poweroff

