# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_10M.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=10M

