# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_14msec.target
[Unit]
Description=Foo
[Login]
IdleActionSec=14msec

