# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_7seconds.target
[Unit]
Description=Foo
[Login]
IdleActionSec=7seconds

