# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_9min.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=9min

