# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleHibernateKey_poweroff.target
[Unit]
Description=Foo
[Login]
HandleHibernateKey=poweroff

