# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_18year.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=18year

