# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleSuspendKey_halt.target
[Unit]
Description=Foo
[Login]
HandleSuspendKey=halt

