# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_13minute.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=13minute

