# systemdlint/tests/auto/InvalidValue/2.31/good/Login_SuspendKeyIgnoreInhibited_on.target
[Unit]
Description=Foo
[Login]
SuspendKeyIgnoreInhibited=on

