# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_2d.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=2d

