# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_15min.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=15min

