# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleHibernateKey_suspend.target
[Unit]
Description=Foo
[Login]
HandleHibernateKey=suspend

