# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_17minutes.target
[Unit]
Description=Foo
[Login]
IdleActionSec=17minutes

