# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_5hour.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=5hour

