# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_8minute.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=8minute

