# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_26year.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=26year

