# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_17second.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=17second

