# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_6days.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=6days

