# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HibernateKeyIgnoreInhibited_no.target
[Unit]
Description=Foo
[Login]
HibernateKeyIgnoreInhibited=no

