# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_8weeks.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=8weeks

