# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_2us.target
[Unit]
Description=Foo
[Login]
IdleActionSec=2us

