# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_20usec.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=20usec

