# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_4ms.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=4ms

