# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_22hr.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=22hr

