# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_18year.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=18year

