# systemdlint/tests/auto/InvalidValue/2.31/good/Login_LidSwitchIgnoreInhibited_1.target
[Unit]
Description=Foo
[Login]
LidSwitchIgnoreInhibited=1

