# systemdlint/tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_9month.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=9month

