# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_4ms.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=4ms

