# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_4hour.target
[Unit]
Description=Foo
[Login]
IdleActionSec=4hour

