# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_7seconds.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=7seconds

