# systemdlint/tests/auto/InvalidValue/2.31/good/Login_RuntimeDirectorySize_7168.target
[Unit]
Description=Foo
[Login]
RuntimeDirectorySize=7168

