# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitch_reboot.target
[Unit]
Description=Foo
[Login]
HandleLidSwitch=reboot

