# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_11m.target
[Unit]
Description=Foo
[Login]
IdleActionSec=11m

