# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandlePowerKey_hibernate.target
[Unit]
Description=Foo
[Login]
HandlePowerKey=hibernate

