# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_22hr.target
[Unit]
Description=Foo
[Login]
IdleActionSec=22hr

