# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_19d.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=19d

