# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitch_hybrid-sleep.target
[Unit]
Description=Foo
[Login]
HandleLidSwitch=hybrid-sleep

