# systemdlint/tests/auto/InvalidValue/2.31/good/Login_RuntimeDirectorySize_1024.target
[Unit]
Description=Foo
[Login]
RuntimeDirectorySize=1024

