# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_23m.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=23m

