# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_4hour.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=4hour

