# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_28%C2%B5s.target
[Unit]
Description=Foo
[Login]
IdleActionSec=28µs

