# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_28h.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=28h

