# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_15ms.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=15ms

