# systemdlint/tests/auto/InvalidValue/2.31/good/Login_RuntimeDirectorySize_2048.target
[Unit]
Description=Foo
[Login]
RuntimeDirectorySize=2048

