# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_16week.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=16week

