# systemdlint/tests/auto/InvalidValue/2.31/good/Login_KillUserProcesses_no.target
[Unit]
Description=Foo
[Login]
KillUserProcesses=no

