# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_8minute.target
[Unit]
Description=Foo
[Login]
IdleActionSec=8minute

