# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_20usec.target
[Unit]
Description=Foo
[Login]
IdleActionSec=20usec

