# systemdlint/tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_24years.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=24years

