# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_5hour.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=5hour

