# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_0sec.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=0sec

