# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_12second.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=12second

