# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_28h.target
[Unit]
Description=Foo
[Login]
IdleActionSec=28h

