# systemdlint/tests/auto/InvalidValue/2.31/good/Login_SuspendKeyIgnoreInhibited_false.target
[Unit]
Description=Foo
[Login]
SuspendKeyIgnoreInhibited=false

