# systemdlint/tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_9month.target
[Unit]
Description=Foo
[Login]
IdleActionSec=9month

