# systemdlint/tests/auto/InvalidValue/2.31/good/Login_HandleHibernateKey_hybrid-sleep.target
[Unit]
Description=Foo
[Login]
HandleHibernateKey=hybrid-sleep

