# systemdlint/tests/auto/InvalidValue/2.31/good/MACVTAP_Mode_source.target
[Unit]
Description=Foo
[MACVTAP]
Mode=source

