# systemdlint/tests/auto/InvalidValue/2.31/good/MACVLAN_Mode_source.target
[Unit]
Description=Foo
[MACVLAN]
Mode=source

