# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_16seconds.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=16seconds

