# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultCPUAccounting_true.target
[Unit]
Description=Foo
[Manager]
DefaultCPUAccounting=true

