# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_20w.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=20w

