# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_11m.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=11m

