# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultIOAccounting_on.target
[Unit]
Description=Foo
[Manager]
DefaultIOAccounting=on

