# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultCPUAccounting_false.target
[Unit]
Description=Foo
[Manager]
DefaultCPUAccounting=false

