# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_28h.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=28h

