# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_22weeks.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=22weeks

