# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_22weeks.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=22weeks

