# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_1day.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=1day

