# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_0.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=0

