# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_13minute.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=13minute

