# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_15min.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=15min

