# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_27y.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=27y

