# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStartSec_8minute.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStartSec=8minute

