# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_18sec.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=18sec

