# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_6h.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=6h

