# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultLimitCORE_abc.target
[Unit]
Description=Foo
[Manager]
DefaultLimitCORE=abc

