# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_20usec.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=20usec

