# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultLimitRSS_abc.target
[Unit]
Description=Foo
[Manager]
DefaultLimitRSS=abc

