# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_18sec.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=18sec

