# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_LogColor_on.target
[Unit]
Description=Foo
[Manager]
LogColor=on

