# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_1day.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=1day

