# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_CrashReboot_no.target
[Unit]
Description=Foo
[Manager]
CrashReboot=no

