# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_25years.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=25years

