# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultIPAccounting_0.target
[Unit]
Description=Foo
[Manager]
DefaultIPAccounting=0

