# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_15ms.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=15ms

