# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_23week.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=23week

