# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_LogColor_false.target
[Unit]
Description=Foo
[Manager]
LogColor=false

