# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_CrashShell_on.target
[Unit]
Description=Foo
[Manager]
CrashShell=on

