# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_8weeks.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=8weeks

