# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_10M.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=10M

