# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultIOAccounting_false.target
[Unit]
Description=Foo
[Manager]
DefaultIOAccounting=false

