# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_20usec.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=20usec

