# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_21usec.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=21usec

