# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_16week.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=16week

