# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_24years.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=24years

