# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_8minute.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=8minute

