# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_25years.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=25years

