# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_0.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=0

