# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTasksAccounting_true.target
[Unit]
Description=Foo
[Manager]
DefaultTasksAccounting=true

