# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultIPAccounting_off.target
[Unit]
Description=Foo
[Manager]
DefaultIPAccounting=off

