# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_6h.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=6h

