# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultLimitCPU_abc.target
[Unit]
Description=Foo
[Manager]
DefaultLimitCPU=abc

