# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_CPUAffinity_9999999999999.target
[Unit]
Description=Foo
[Manager]
CPUAffinity=9999999999999

