# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultIPAccounting_true.target
[Unit]
Description=Foo
[Manager]
DefaultIPAccounting=true

