# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShowStatus_true.target
[Unit]
Description=Foo
[Manager]
ShowStatus=true

