# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DumpCore_on.target
[Unit]
Description=Foo
[Manager]
DumpCore=on

