# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_27msec.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=27msec

