# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_12month.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=12month

