# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_7seconds.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=7seconds

