# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_15ms.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=15ms

