# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_22hr.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=22hr

