# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultIPAccounting_no.target
[Unit]
Description=Foo
[Manager]
DefaultIPAccounting=no

