# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_1%C2%B5s.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=1µs

