# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_25years.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=25years

