# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_infinity.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=infinity

