# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_10M.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=10M

