# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_9min.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=9min

