# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_0days.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=0days

