# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_CrashChangeVT_on.target
[Unit]
Description=Foo
[Manager]
CrashChangeVT=on

