# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_7minutes.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=7minutes

