# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_10M.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=10M

