# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_7seconds.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=7seconds

