# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_23week.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=23week

