# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_8weeks.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=8weeks

