# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DumpCore_no.target
[Unit]
Description=Foo
[Manager]
DumpCore=no

