# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTasksAccounting_false.target
[Unit]
Description=Foo
[Manager]
DefaultTasksAccounting=false

