# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultIOAccounting_off.target
[Unit]
Description=Foo
[Manager]
DefaultIOAccounting=off

