# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultLimitNOFILE_abc.target
[Unit]
Description=Foo
[Manager]
DefaultLimitNOFILE=abc

