# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_0days.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=0days

