# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_0sec.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=0sec

