# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_28h.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=28h

