# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_infinity.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=infinity

