# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_6days.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=6days

