# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_28%C2%B5s.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=28µs

