# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_16seconds.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=16seconds

