# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_23week.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=23week

