# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_24years.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=24years

