# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultLimitSIGPENDING_abc.target
[Unit]
Description=Foo
[Manager]
DefaultLimitSIGPENDING=abc

