# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_5hr.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=5hr

