# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTasksMax_0%25.target
[Unit]
Description=Foo
[Manager]
DefaultTasksMax=0%

