# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_15ms.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=15ms

