# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_18sec.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=18sec

