# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStartSec_21us.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStartSec=21us

