# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_25s.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=25s

