# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_19d.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=19d

