# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DumpCore_off.target
[Unit]
Description=Foo
[Manager]
DumpCore=off

