# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DumpCore_0.target
[Unit]
Description=Foo
[Manager]
DumpCore=0

