# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_2us.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=2us

