# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShowStatus_0.target
[Unit]
Description=Foo
[Manager]
ShowStatus=0

