# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_9999999999999.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=9999999999999

