# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTasksAccounting_yes.target
[Unit]
Description=Foo
[Manager]
DefaultTasksAccounting=yes

