# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_28h.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=28h

