# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_CPUAffinity_0.target
[Unit]
Description=Foo
[Manager]
CPUAffinity=0

