# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_17second.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=17second

