# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_27y.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=27y

