# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_12month.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=12month

