# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_21usec.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=21usec

