# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultBlockIOAccounting_true.target
[Unit]
Description=Foo
[Manager]
DefaultBlockIOAccounting=true

