# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_3hours.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=3hours

