# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_24w.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=24w

