# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_21us.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=21us

