# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_20usec.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=20usec

