# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_17minutes.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=17minutes

