# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_21us.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=21us

