# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_8minute.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=8minute

