# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStartSec_infinity.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStartSec=infinity

