# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShowStatus_false.target
[Unit]
Description=Foo
[Manager]
ShowStatus=false

