# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_LogColor_1.target
[Unit]
Description=Foo
[Manager]
LogColor=1

