# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_0days.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=0days

