# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTasksAccounting_0.target
[Unit]
Description=Foo
[Manager]
DefaultTasksAccounting=0

