# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_12second.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=12second

