# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_10M.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=10M

