# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_19s.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=19s

