# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_17second.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=17second

