# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_9999999999999.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=9999999999999

