# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_9999999999999.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=9999999999999

