# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStartSec_17second.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStartSec=17second

