# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_21usec.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=21usec

