# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_14M.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=14M

