# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultIPAccounting_yes.target
[Unit]
Description=Foo
[Manager]
DefaultIPAccounting=yes

