# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_26months.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=26months

