# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_infinity.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=infinity

