# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimerAccuracySec_9min.target
[Unit]
Description=Foo
[Manager]
DefaultTimerAccuracySec=9min

