# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultCPUAccounting_0.target
[Unit]
Description=Foo
[Manager]
DefaultCPUAccounting=0

