# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_28%C2%B5s.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=28µs

