# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultLimitMEMLOCK_abc.target
[Unit]
Description=Foo
[Manager]
DefaultLimitMEMLOCK=abc

