# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_9999999999999.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=9999999999999

