# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_25years.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=25years

