# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_WatchdogDevice_--dev--bar.target
[Unit]
Description=Foo
[Manager]
WatchdogDevice=/dev/bar

