# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_13months.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=13months

