# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_27y.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=27y

