# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_6h.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=6h

