# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimerAccuracySec_26year.target
[Unit]
Description=Foo
[Manager]
DefaultTimerAccuracySec=26year

