# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_11m.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=11m

