# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultMemoryAccounting_on.target
[Unit]
Description=Foo
[Manager]
DefaultMemoryAccounting=on

