# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_14msec.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=14msec

