# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_11m.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=11m

