# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_2d.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=2d

