# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultLimitRTTIME_abc.target
[Unit]
Description=Foo
[Manager]
DefaultLimitRTTIME=abc

