# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_27y.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=27y

