# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_28%C2%B5s.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=28µs

