# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_22weeks.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=22weeks

