# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStartSec_24w.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStartSec=24w

