# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStartSec_6h.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStartSec=6h

