# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultCPUAccounting_yes.target
[Unit]
Description=Foo
[Manager]
DefaultCPUAccounting=yes

