# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_12month.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=12month

