# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultLimitRTPRIO_abc.target
[Unit]
Description=Foo
[Manager]
DefaultLimitRTPRIO=abc

