# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_infinity.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=infinity

