# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_19d.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=19d

