# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimerAccuracySec_10M.target
[Unit]
Description=Foo
[Manager]
DefaultTimerAccuracySec=10M

