# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_9min.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=9min

