# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_26year.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=26year

