# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_4hour.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=4hour

