# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_2d.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=2d

