# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimerAccuracySec_16seconds.target
[Unit]
Description=Foo
[Manager]
DefaultTimerAccuracySec=16seconds

