# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_12month.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=12month

