# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_28%C2%B5s.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=28µs

