# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimerAccuracySec_12month.target
[Unit]
Description=Foo
[Manager]
DefaultTimerAccuracySec=12month

