# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_18sec.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=18sec

