# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_12month.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=12month

