# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultLimitNICE_abc.target
[Unit]
Description=Foo
[Manager]
DefaultLimitNICE=abc

