# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_22weeks.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=22weeks

