# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_8minute.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=8minute

