# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_9999999999999.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=9999999999999

