# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_LogColor_true.target
[Unit]
Description=Foo
[Manager]
LogColor=true

