# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_11y.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=11y

