# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_24w.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=24w

