# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_9min.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=9min

