# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_9month.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=9month

