# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_16week.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=16week

