# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_19s.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=19s

