# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_18year.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=18year

