# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_4hour.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=4hour

