# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_24w.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=24w

