# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_16seconds.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=16seconds

