# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultLimitLOCKS_abc.target
[Unit]
Description=Foo
[Manager]
DefaultLimitLOCKS=abc

