# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_18sec.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=18sec

