# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_26months.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=26months

