# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_7minutes.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=7minutes

