# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_22hr.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=22hr

