# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_2d.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=2d

