# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_23m.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=23m

