# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_15ms.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=15ms

