# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultMemoryAccounting_off.target
[Unit]
Description=Foo
[Manager]
DefaultMemoryAccounting=off

