# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_10M.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=10M

