# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_LogLevel_abc.target
[Unit]
Description=Foo
[Manager]
LogLevel=abc

