# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_18year.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=18year

