# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultCPUAccounting_off.target
[Unit]
Description=Foo
[Manager]
DefaultCPUAccounting=off

