# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_28%C2%B5s.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=28µs

