# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_4hour.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=4hour

