# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_10M.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=10M

