# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_8minute.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=8minute

