# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_18sec.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=18sec

