# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_15min.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=15min

