# systemdlint/tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_1%C2%B5s.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=1µs

