# systemdlint/tests/auto/InvalidValue/2.31/good/Match_Type_abc.target
[Unit]
Description=Foo
[Match]
Type=abc

