# systemdlint/tests/auto/InvalidValue/2.31/good/Match_Virtualization_abc.target
[Unit]
Description=Foo
[Match]
Virtualization=abc

