# systemdlint/tests/auto/InvalidValue/2.31/good/Match_KernelCommandLine_abc.target
[Unit]
Description=Foo
[Match]
KernelCommandLine=abc

