# systemdlint/tests/auto/InvalidValue/2.31/good/Match_KernelVersion_abc.target
[Unit]
Description=Foo
[Match]
KernelVersion=abc

