# systemdlint/tests/auto/InvalidValue/2.31/good/Match_MACAddress_abc.target
[Unit]
Description=Foo
[Match]
MACAddress=abc

