# systemdlint/tests/auto/InvalidValue/2.31/good/Match_Architecture_abc.target
[Unit]
Description=Foo
[Match]
Architecture=abc

