# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_LimitCPU_0.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
LimitCPU=0

