# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_UtmpIdentifier_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
UtmpIdentifier=abc

