# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_LimitRSS_9999999999999.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
LimitRSS=9999999999999

