# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_MemoryLimit_0M.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
MemoryLimit=0M

