# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_TTYReset_1.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
TTYReset=1

