# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_ProtectHome_yes.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
ProtectHome=yes

