# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_CPUAffinity_numa.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
CPUAffinity=numa

