# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_MemoryLow_0B.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
MemoryLow=0B

