# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_IOWeight_10000.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
IOWeight=10000

