# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_TTYVTDisallocate_false.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
TTYVTDisallocate=false

