# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_Personality_s390x.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
Personality=s390x

