# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_MemoryHigh_0T.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
MemoryHigh=0T

