# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_RemoveIPC_off.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
RemoveIPC=off

