# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_UMask_9999999999999.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
UMask=9999999999999

