# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_RemoveIPC_1.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
RemoveIPC=1

