# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_BlockIODeviceWeight_--dev--foo%2010000.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
BlockIODeviceWeight=/dev/foo 10000

