# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_SystemCallErrorNumber_EL3RST.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
SystemCallErrorNumber=EL3RST

