# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_TimeoutSec_27msec.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
TimeoutSec=27msec

