# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_SyslogFacility_local2.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
SyslogFacility=local2

