# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_TasksAccounting_on.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
TasksAccounting=on

