# systemdlint/tests/auto/InvalidValue/2.31/good/Mount_BlockStartupIOWeight_1.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
BlockStartupIOWeight=1

