# systemdlint/tests/auto/InvalidValue/2.31/good/Neighbor_MACAddress_abc.target
[Unit]
Description=Foo
[Neighbor]
MACAddress=abc

