# systemdlint/tests/auto/InvalidValue/2.31/good/NetDev_Kind_ipvlan.target
[Unit]
Description=Foo
[NetDev]
Name=abc
Kind=ipvlan

