# systemdlint/tests/auto/InvalidValue/2.31/good/NetDev_Kind_l2tp.target
[Unit]
Description=Foo
[NetDev]
Name=abc
Kind=l2tp

