# systemdlint/tests/auto/InvalidValue/2.31/good/NetDev_Kind_ip6gretap.target
[Unit]
Description=Foo
[NetDev]
Name=abc
Kind=ip6gretap

