# systemdlint/tests/auto/InvalidValue/2.31/good/NetDev_MACAddress_abc.target
[Unit]
Description=Foo
[NetDev]
Name=abc
MACAddress=abc

