# systemdlint/tests/auto/InvalidValue/2.31/good/Network_ConfigureWithoutCarrier_no.target
[Unit]
Description=Foo
[Network]
ConfigureWithoutCarrier=no

