# systemdlint/tests/auto/InvalidValue/2.31/good/Network_IPVLAN_abc.target
[Unit]
Description=Foo
[Network]
IPVLAN=abc

