# systemdlint/tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_18sec.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=18sec

