# systemdlint/tests/auto/InvalidValue/2.31/good/Network_ActiveSlave_false.target
[Unit]
Description=Foo
[Network]
ActiveSlave=false

