# systemdlint/tests/auto/InvalidValue/2.31/good/Network_MulticastDNS_yes.target
[Unit]
Description=Foo
[Network]
MulticastDNS=yes

