# systemdlint/tests/auto/InvalidValue/2.31/good/Network_NTP_http%3A----user%3Apwd%40test.com%3A1234.target
[Unit]
Description=Foo
[Network]
NTP=http://user:pwd@test.com:1234

