# systemdlint/tests/auto/InvalidValue/2.31/good/Network_MulticastDNS_0.target
[Unit]
Description=Foo
[Network]
MulticastDNS=0

