# systemdlint/tests/auto/InvalidValue/2.31/good/Network_LLDP_0.target
[Unit]
Description=Foo
[Network]
LLDP=0

