# systemdlint/tests/auto/InvalidValue/2.31/good/Network_DNSOverTLS_opportunistic.target
[Unit]
Description=Foo
[Network]
DNSOverTLS=opportunistic

