# systemdlint/tests/auto/InvalidValue/2.31/good/Network_LLMNR_resolve.target
[Unit]
Description=Foo
[Network]
LLMNR=resolve

