# systemdlint/tests/auto/InvalidValue/2.31/good/Network_DHCP_no.target
[Unit]
Description=Foo
[Network]
DHCP=no

