# systemdlint/tests/auto/InvalidValue/2.31/good/Network_IPMasquerade_on.target
[Unit]
Description=Foo
[Network]
IPMasquerade=on

