# systemdlint/tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_17second.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=17second

