# systemdlint/tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_on.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=on

