# systemdlint/tests/auto/InvalidValue/2.31/good/Network_LLMNR_false.target
[Unit]
Description=Foo
[Network]
LLMNR=false

