# systemdlint/tests/auto/InvalidValue/2.31/good/Network_ActiveSlave_0.target
[Unit]
Description=Foo
[Network]
ActiveSlave=0

