# systemdlint/tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_infinity.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=infinity

