# systemdlint/tests/auto/InvalidValue/2.31/good/Network_Bond_abc.target
[Unit]
Description=Foo
[Network]
Bond=abc

