# systemdlint/tests/auto/InvalidValue/2.31/good/Network_DHCP_ipv4.target
[Unit]
Description=Foo
[Network]
DHCP=ipv4

