# systemdlint/tests/auto/InvalidValue/2.31/good/Network_NTP_file%3A----test.a.b.com.target
[Unit]
Description=Foo
[Network]
NTP=file://test.a.b.com

