# systemdlint/tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_21us.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=21us

