# systemdlint/tests/auto/InvalidValue/2.31/good/Network_MulticastDNS_false.target
[Unit]
Description=Foo
[Network]
MulticastDNS=false

