# systemdlint/tests/auto/InvalidValue/2.31/good/Network_LinkLocalAddressing_yes.target
[Unit]
Description=Foo
[Network]
LinkLocalAddressing=yes

