# systemdlint/tests/auto/InvalidValue/2.31/good/Network_DNSSEC_allow-downgrade.target
[Unit]
Description=Foo
[Network]
DNSSEC=allow-downgrade

