# systemdlint/tests/auto/InvalidValue/2.31/good/Network_IPv4LLRoute_false.target
[Unit]
Description=Foo
[Network]
IPv4LLRoute=false

