# systemdlint/tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_23week.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=23week

