# systemdlint/tests/auto/InvalidValue/2.31/good/Network_ConfigureWithoutCarrier_1.target
[Unit]
Description=Foo
[Network]
ConfigureWithoutCarrier=1

