# systemdlint/tests/auto/InvalidValue/2.31/good/Network_EmitLLDP_nearest-bridge.target
[Unit]
Description=Foo
[Network]
EmitLLDP=nearest-bridge

