# systemdlint/tests/auto/InvalidValue/2.31/good/Network_DHCP_yes.target
[Unit]
Description=Foo
[Network]
DHCP=yes

