# systemdlint/tests/auto/InvalidValue/2.31/good/Network_LLDP_routers-only.target
[Unit]
Description=Foo
[Network]
LLDP=routers-only

