# systemdlint/tests/auto/InvalidValue/2.31/good/Network_NTP_http%3A----a.b.com.target
[Unit]
Description=Foo
[Network]
NTP=http://a.b.com

