# systemdlint/tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_6h.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=6h

