# systemdlint/tests/auto/InvalidValue/2.31/good/Network_MulticastDNS_resolve.target
[Unit]
Description=Foo
[Network]
MulticastDNS=resolve

