# systemdlint/tests/auto/InvalidValue/2.31/good/Network_DNSDefaultRoute_false.target
[Unit]
Description=Foo
[Network]
DNSDefaultRoute=false

