# systemdlint/tests/auto/InvalidValue/2.31/good/Network_DNSDefaultRoute_off.target
[Unit]
Description=Foo
[Network]
DNSDefaultRoute=off

