# systemdlint/tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_4hour.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=4hour

