# systemdlint/tests/auto/InvalidValue/2.31/good/Network_IPMasquerade_no.target
[Unit]
Description=Foo
[Network]
IPMasquerade=no

