# systemdlint/tests/auto/InvalidValue/2.31/good/Network_DHCPServer_no.target
[Unit]
Description=Foo
[Network]
DHCPServer=no

