# systemdlint/tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_0.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=0

