# systemdlint/tests/auto/InvalidValue/2.31/good/Network_IPMasquerade_yes.target
[Unit]
Description=Foo
[Network]
IPMasquerade=yes

