# systemdlint/tests/auto/InvalidValue/2.31/good/Network_IPMasquerade_0.target
[Unit]
Description=Foo
[Network]
IPMasquerade=0

