# systemdlint/tests/auto/InvalidValue/2.31/good/Network_MulticastDNS_no.target
[Unit]
Description=Foo
[Network]
MulticastDNS=no

