# systemdlint/tests/auto/InvalidValue/2.31/good/Network_DNSSEC_true.target
[Unit]
Description=Foo
[Network]
DNSSEC=true

