# systemdlint/tests/auto/InvalidValue/2.31/good/Network_IPv4ProxyARP_0.target
[Unit]
Description=Foo
[Network]
IPv4ProxyARP=0

