# systemdlint/tests/auto/InvalidValue/2.31/good/Network_IPv6PrefixDelegation_yes.target
[Unit]
Description=Foo
[Network]
IPv6PrefixDelegation=yes

