# systemdlint/tests/auto/InvalidValue/2.31/good/Network_EmitLLDP_customer-bridge.target
[Unit]
Description=Foo
[Network]
EmitLLDP=customer-bridge

