# systemdlint/tests/auto/InvalidValue/2.31/good/Network_DNSSEC_false.target
[Unit]
Description=Foo
[Network]
DNSSEC=false

