# systemdlint/tests/auto/InvalidValue/2.31/good/Network_LLDP_false.target
[Unit]
Description=Foo
[Network]
LLDP=false

