# systemdlint/tests/auto/InvalidValue/2.31/good/Network_ConfigureWithoutCarrier_on.target
[Unit]
Description=Foo
[Network]
ConfigureWithoutCarrier=on

