# systemdlint/tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_8minute.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=8minute

