# systemdlint/tests/auto/InvalidValue/2.31/good/Path_PathExists_--dev--.target
[Unit]
Description=Foo
[Path]
PathExists=/dev/

