# systemdlint/tests/auto/InvalidValue/2.31/good/Path_PathModified_--dev--.target
[Unit]
Description=Foo
[Path]
PathModified=/dev/

