# systemdlint/tests/auto/InvalidValue/2.31/good/Path_MakeDirectory_off.target
[Unit]
Description=Foo
[Path]
MakeDirectory=off

