# systemdlint/tests/auto/InvalidValue/2.31/good/Path_MakeDirectory_on.target
[Unit]
Description=Foo
[Path]
MakeDirectory=on

