# systemdlint/tests/auto/InvalidValue/2.31/good/Path_MakeDirectory_true.target
[Unit]
Description=Foo
[Path]
MakeDirectory=true

