# systemdlint/tests/auto/InvalidValue/2.31/good/Path_DirectoryMode_007.target
[Unit]
Description=Foo
[Path]
DirectoryMode=007

