# systemdlint/tests/auto/InvalidValue/2.31/good/Peer_Name_abc.target
[Unit]
Description=Foo
[Peer]
Name=abc
Name=abc

