# systemdlint/tests/auto/InvalidValue/2.31/good/Peer_MACAddress_abc.target
[Unit]
Description=Foo
[Peer]
Name=abc
MACAddress=abc

