# systemdlint/tests/auto/InvalidValue/2.31/good/Route_Gateway__dhcp.target
[Unit]
Description=Foo
[Route]
Gateway=_dhcp

