# systemdlint/tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_OutgoingInterface_abc.target
[Unit]
Description=Foo
[RoutingPolicyRule]
OutgoingInterface=abc

