# systemdlint/tests/auto/InvalidValue/2.31/good/Route_Table_1.target
[Unit]
Description=Foo
[Route]
Table=1

