# systemdlint/tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_SourcePort_9999999999999%209999999999999.target
[Unit]
Description=Foo
[RoutingPolicyRule]
SourcePort=9999999999999 9999999999999

