# systemdlint/tests/auto/InvalidValue/2.31/good/Route_MTUBytes_7168.target
[Unit]
Description=Foo
[Route]
MTUBytes=7168

