# systemdlint/tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_InvertRule_true.target
[Unit]
Description=Foo
[RoutingPolicyRule]
InvertRule=true

