# systemdlint/tests/auto/InvalidValue/2.31/good/Route_Table_4294967295.target
[Unit]
Description=Foo
[Route]
Table=4294967295

