# systemdlint/tests/auto/InvalidValue/2.31/good/Route_InitialAdvertisedReceiveWindow_1.target
[Unit]
Description=Foo
[Route]
InitialAdvertisedReceiveWindow=1

