# systemdlint/tests/auto/InvalidValue/2.31/good/Route_InitialCongestionWindow_1.target
[Unit]
Description=Foo
[Route]
InitialCongestionWindow=1

