# systemdlint/tests/auto/InvalidValue/2.31/good/Route_QuickAck_yes.target
[Unit]
Description=Foo
[Route]
QuickAck=yes

