# systemdlint/tests/auto/InvalidValue/2.31/good/Route_Type_unreachable.target
[Unit]
Description=Foo
[Route]
Type=unreachable

