# systemdlint/tests/auto/InvalidValue/2.31/good/Route_Protocol_static.target
[Unit]
Description=Foo
[Route]
Protocol=static

