# systemdlint/tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_TypeOfService_255.target
[Unit]
Description=Foo
[RoutingPolicyRule]
TypeOfService=255

