# systemdlint/tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_Priority_9999999999999.target
[Unit]
Description=Foo
[RoutingPolicyRule]
Priority=9999999999999

