# systemdlint/tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_InvertRule_yes.target
[Unit]
Description=Foo
[RoutingPolicyRule]
InvertRule=yes

