# systemdlint/tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_IPProtocol_udp.target
[Unit]
Description=Foo
[RoutingPolicyRule]
IPProtocol=udp

