# systemdlint/tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_FirewallMark_4294967295.target
[Unit]
Description=Foo
[RoutingPolicyRule]
FirewallMark=4294967295

