# systemdlint/tests/auto/InvalidValue/2.31/good/Route_GatewayOnLink_1.target
[Unit]
Description=Foo
[Route]
GatewayOnLink=1

