# systemdlint/tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_SourcePort_0%200.target
[Unit]
Description=Foo
[RoutingPolicyRule]
SourcePort=0 0

