# systemdlint/tests/auto/InvalidValue/2.31/good/Route_Scope_host.target
[Unit]
Description=Foo
[Route]
Scope=host

