# systemdlint/tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_IncomingInterface_abc.target
[Unit]
Description=Foo
[RoutingPolicyRule]
IncomingInterface=abc

