# systemdlint/tests/auto/InvalidValue/2.31/good/Route_GatewayOnLink_no.target
[Unit]
Description=Foo
[Route]
GatewayOnLink=no

