# systemdlint/tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_SourcePort_9999999999999%200.target
[Unit]
Description=Foo
[RoutingPolicyRule]
SourcePort=9999999999999 0

