# systemdlint/tests/auto/InvalidValue/2.31/good/Route_Type_prohibit.target
[Unit]
Description=Foo
[Route]
Type=prohibit

