# systemdlint/tests/auto/InvalidValue/2.31/good/Route_QuickAck_true.target
[Unit]
Description=Foo
[Route]
QuickAck=true

