# systemdlint/tests/auto/InvalidValue/2.31/good/Route_GatewayOnLink_on.target
[Unit]
Description=Foo
[Route]
GatewayOnLink=on

