# systemdlint/tests/auto/InvalidValue/2.31/good/Route_GatewayOnLink_off.target
[Unit]
Description=Foo
[Route]
GatewayOnLink=off

