# systemdlint/tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_DestinationPort_9999999999999%200.target
[Unit]
Description=Foo
[RoutingPolicyRule]
DestinationPort=9999999999999 0

