# systemdlint/tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_InvertRule_on.target
[Unit]
Description=Foo
[RoutingPolicyRule]
InvertRule=on

