# systemdlint/tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_IPProtocol_17.target
[Unit]
Description=Foo
[RoutingPolicyRule]
IPProtocol=17

