# systemdlint/tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_SourcePort_0%209999999999999.target
[Unit]
Description=Foo
[RoutingPolicyRule]
SourcePort=0 9999999999999

