# systemdlint/tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_IPProtocol_sctp.target
[Unit]
Description=Foo
[RoutingPolicyRule]
IPProtocol=sctp

