# systemdlint/tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_InvertRule_1.target
[Unit]
Description=Foo
[RoutingPolicyRule]
InvertRule=1

