# systemdlint/tests/auto/InvalidValue/2.31/good/Route_Scope_link.target
[Unit]
Description=Foo
[Route]
Scope=link

