# systemdlint/tests/auto/InvalidValue/2.31/good/Route_Type_blackhole.target
[Unit]
Description=Foo
[Route]
Type=blackhole

