# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitCORE_0.target
[Unit]
Description=Foo
[Service]
LimitCORE=0

