# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_13minute.target
[Unit]
Description=Foo
[Service]
RestartSec=13minute

