# systemdlint/tests/auto/InvalidValue/2.31/good/Service_BlockIODeviceWeight_--dev--bar%2010000.target
[Unit]
Description=Foo
[Service]
BlockIODeviceWeight=/dev/bar 10000

