# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SendSIGHUP_true.target
[Unit]
Description=Foo
[Service]
SendSIGHUP=true

