# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_12second.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=12second

