# systemdlint/tests/auto/InvalidValue/2.31/good/Service_MountAPIVFS_yes.target
[Unit]
Description=Foo
[Service]
MountAPIVFS=yes

