# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EISCONN.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EISCONN

