# systemdlint/tests/auto/InvalidValue/2.31/good/Service_MemoryDenyWriteExecute_on.target
[Unit]
Description=Foo
[Service]
MemoryDenyWriteExecute=on

