# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitRSS_0.target
[Unit]
Description=Foo
[Service]
LimitRSS=0

