# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_19d.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=19d

