# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGHUP.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGHUP

