# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartLimitInterval_11y.target
[Unit]
Description=Foo
[Service]
StartLimitInterval=11y

