# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartLimitInterval_13minute.target
[Unit]
Description=Foo
[Service]
StartLimitInterval=13minute

