# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_23week.target
[Unit]
Description=Foo
[Service]
WatchdogSec=23week

