# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ECANCELED.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ECANCELED

