# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_13minute.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=13minute

