# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_1day.target
[Unit]
Description=Foo
[Service]
RestartSec=1day

