# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_18year.target
[Unit]
Description=Foo
[Service]
RestartSec=18year

