# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_4hour.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=4hour

