# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EINVAL.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EINVAL

