# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ECONNRESET.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ECONNRESET

