# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGCONT.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGCONT

