# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_15ms.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=15ms

