# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_0days.target
[Unit]
Description=Foo
[Service]
RestartSec=0days

