# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_9999999999999.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=9999999999999

