# systemdlint/tests/auto/InvalidValue/2.31/good/Service_ProtectHome_tmpfs.target
[Unit]
Description=Foo
[Service]
ProtectHome=tmpfs

