# systemdlint/tests/auto/InvalidValue/2.31/good/Service_Group_metoo.target
[Unit]
Description=Foo
[Service]
Group=metoo

