# systemdlint/tests/auto/InvalidValue/2.31/good/Service_BlockStartupIOWeight_1.target
[Unit]
Description=Foo
[Service]
BlockStartupIOWeight=1

