# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_17minutes.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=17minutes

