# systemdlint/tests/auto/InvalidValue/2.31/good/Service_ExecStopPost_--bin--sh.target
[Unit]
Description=Foo
[Service]
ExecStopPost=/bin/sh

