# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGINFO.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGINFO

