# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SyslogLevel_crit.target
[Unit]
Description=Foo
[Service]
SyslogLevel=crit

