# systemdlint/tests/auto/InvalidValue/2.31/good/Service_BlockIOWeight_10000.target
[Unit]
Description=Foo
[Service]
BlockIOWeight=10000

