# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGIO.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGIO

