# systemdlint/tests/auto/InvalidValue/2.31/good/Service_PrivateNetwork_true.target
[Unit]
Description=Foo
[Service]
PrivateNetwork=true

