# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_13.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=13

