# systemdlint/tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGCONT.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGCONT

