# systemdlint/tests/auto/InvalidValue/2.31/good/Service_IOAccounting_off.target
[Unit]
Description=Foo
[Service]
IOAccounting=off

