# systemdlint/tests/auto/InvalidValue/2.31/good/Service_IOAccounting_0.target
[Unit]
Description=Foo
[Service]
IOAccounting=0

