# systemdlint/tests/auto/InvalidValue/2.31/good/Service_Restart_on-abnormal.target
[Unit]
Description=Foo
[Service]
Restart=on-abnormal

