# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitRTPRIO_0.target
[Unit]
Description=Foo
[Service]
LimitRTPRIO=0

