# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_DATAERR.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=DATAERR

