# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_infinity.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=infinity

