# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGILL.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGILL

