# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_11.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=11

