# systemdlint/tests/auto/InvalidValue/2.31/good/Service_ProtectHome_on.target
[Unit]
Description=Foo
[Service]
ProtectHome=on

