# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RootDirectoryStartOnly_true.target
[Unit]
Description=Foo
[Service]
RootDirectoryStartOnly=true

