# systemdlint/tests/auto/InvalidValue/2.31/good/Service_CPUSchedulingResetOnFork_false.target
[Unit]
Description=Foo
[Service]
CPUSchedulingResetOnFork=false

