# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGALRM.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGALRM

