# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_56.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=56

