# systemdlint/tests/auto/InvalidValue/2.31/good/Service_CPUSchedulingResetOnFork_true.target
[Unit]
Description=Foo
[Service]
CPUSchedulingResetOnFork=true

