# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestrictSUIDSGID_true.target
[Unit]
Description=Foo
[Service]
RestrictSUIDSGID=true

