# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RemainAfterExit_0.target
[Unit]
Description=Foo
[Service]
RemainAfterExit=0

