# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGWINCH.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGWINCH

