# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_5hour.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=5hour

