# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGEMT.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGEMT

