# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_2d.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=2d

