# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SendSIGHUP_0.target
[Unit]
Description=Foo
[Service]
SendSIGHUP=0

