# systemdlint/tests/auto/InvalidValue/2.31/good/Service_MemoryAccounting_false.target
[Unit]
Description=Foo
[Service]
MemoryAccounting=false

