# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EINPROGRESS.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EINPROGRESS

