# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartLimitIntervalSec_8minute.target
[Unit]
Description=Foo
[Service]
StartLimitIntervalSec=8minute

