# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestrictSUIDSGID_on.target
[Unit]
Description=Foo
[Service]
RestrictSUIDSGID=on

