# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestrictRealtime_on.target
[Unit]
Description=Foo
[Service]
RestrictRealtime=on

