# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SyslogLevelPrefix_0.target
[Unit]
Description=Foo
[Service]
SyslogLevelPrefix=0

