# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_20.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=20

