# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_23week.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=23week

