# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SyslogLevel_notice.target
[Unit]
Description=Foo
[Service]
SyslogLevel=notice

