# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_7minutes.target
[Unit]
Description=Foo
[Service]
RestartSec=7minutes

