# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_11y.target
[Unit]
Description=Foo
[Service]
TimeoutSec=11y

