# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitCPU_infinity.target
[Unit]
Description=Foo
[Service]
LimitCPU=infinity

