# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_19d.target
[Unit]
Description=Foo
[Service]
WatchdogSec=19d

