# systemdlint/tests/auto/InvalidValue/2.31/good/Service_DisableControllers_cpuacct.target
[Unit]
Description=Foo
[Service]
DisableControllers=cpuacct

