# systemdlint/tests/auto/InvalidValue/2.31/good/Service_DisableControllers_devices.target
[Unit]
Description=Foo
[Service]
DisableControllers=devices

