# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_12.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=12

