# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SyslogFacility_kern.target
[Unit]
Description=Foo
[Service]
SyslogFacility=kern

