# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestrictSUIDSGID_false.target
[Unit]
Description=Foo
[Service]
RestrictSUIDSGID=false

