# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestrictRealtime_true.target
[Unit]
Description=Foo
[Service]
RestrictRealtime=true

