# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitSIGPENDING_0.target
[Unit]
Description=Foo
[Service]
LimitSIGPENDING=0

