# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SyslogLevelPrefix_yes.target
[Unit]
Description=Foo
[Service]
SyslogLevelPrefix=yes

