# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RemoveIPC_yes.target
[Unit]
Description=Foo
[Service]
RemoveIPC=yes

