# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_1day.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=1day

