# systemdlint/tests/auto/InvalidValue/2.31/good/Service_IOSchedulingClass_best-effort.target
[Unit]
Description=Foo
[Service]
IOSchedulingClass=best-effort

