# systemdlint/tests/auto/InvalidValue/2.31/good/Service_MemoryLimit_0G.target
[Unit]
Description=Foo
[Service]
MemoryLimit=0G

