# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestrictRealtime_false.target
[Unit]
Description=Foo
[Service]
RestrictRealtime=false

