# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitAS_infinity.target
[Unit]
Description=Foo
[Service]
LimitAS=infinity

