# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGXCPU.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGXCPU

