# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGVTALRM.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGVTALRM

