# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ENOSPC.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ENOSPC

