# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_27msec.target
[Unit]
Description=Foo
[Service]
WatchdogSec=27msec

