# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGCONT.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGCONT

