# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGTSTP.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGTSTP

