# systemdlint/tests/auto/InvalidValue/2.31/good/Service_MemoryDenyWriteExecute_no.target
[Unit]
Description=Foo
[Service]
MemoryDenyWriteExecute=no

