# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_4hour.target
[Unit]
Description=Foo
[Service]
RestartSec=4hour

