# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_24years.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=24years

