# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitNOFILE_infinity.target
[Unit]
Description=Foo
[Service]
LimitNOFILE=infinity

