# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_13minute.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=13minute

