# systemdlint/tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGALRM.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGALRM

