# systemdlint/tests/auto/InvalidValue/2.31/good/Service_PrivateTmp_yes.target
[Unit]
Description=Foo
[Service]
PrivateTmp=yes

