# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TTYReset_1.target
[Unit]
Description=Foo
[Service]
TTYReset=1

