# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TasksMax_9999999999999.target
[Unit]
Description=Foo
[Service]
TasksMax=9999999999999

