# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_21usec.target
[Unit]
Description=Foo
[Service]
TimeoutSec=21usec

