# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_13months.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=13months

