# systemdlint/tests/auto/InvalidValue/2.31/good/Service_CPUSchedulingResetOnFork_on.target
[Unit]
Description=Foo
[Service]
CPUSchedulingResetOnFork=on

