# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ENOMEM.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ENOMEM

