# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RootDirectoryStartOnly_no.target
[Unit]
Description=Foo
[Service]
RootDirectoryStartOnly=no

