# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_0.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=0

