# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitLOCKS_0.target
[Unit]
Description=Foo
[Service]
LimitLOCKS=0

