# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_14msec.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=14msec

