# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_21usec.target
[Unit]
Description=Foo
[Service]
RestartSec=21usec

