# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartLimitInterval_18sec.target
[Unit]
Description=Foo
[Service]
StartLimitInterval=18sec

