# systemdlint/tests/auto/InvalidValue/2.31/good/Service_UtmpMode_init.target
[Unit]
Description=Foo
[Service]
UtmpMode=init

