# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallFilter_abc.target
[Unit]
Description=Foo
[Service]
SystemCallFilter=abc

