# systemdlint/tests/auto/InvalidValue/2.31/good/Service_IOAccounting_false.target
[Unit]
Description=Foo
[Service]
IOAccounting=false

