# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StandardError_syslog%2Bconsole.target
[Unit]
Description=Foo
[Service]
StandardError=syslog+console

