# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartLimitInterval_17minutes.target
[Unit]
Description=Foo
[Service]
StartLimitInterval=17minutes

