# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitCPU_9999999999999.target
[Unit]
Description=Foo
[Service]
LimitCPU=9999999999999

