# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ETIMEDOUT.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ETIMEDOUT

