# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_57.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=57

