# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RootDirectoryStartOnly_1.target
[Unit]
Description=Foo
[Service]
RootDirectoryStartOnly=1

