# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_22hr.target
[Unit]
Description=Foo
[Service]
TimeoutSec=22hr

