# systemdlint/tests/auto/InvalidValue/2.31/good/Service_IOSchedulingPriority_0.target
[Unit]
Description=Foo
[Service]
IOSchedulingPriority=0

