# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_12second.target
[Unit]
Description=Foo
[Service]
WatchdogSec=12second

