# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TasksAccounting_1.target
[Unit]
Description=Foo
[Service]
TasksAccounting=1

