# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_28h.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=28h

