# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StandardOutput_tty.target
[Unit]
Description=Foo
[Service]
StandardOutput=tty

