# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_7minutes.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=7minutes

