# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_22hr.target
[Unit]
Description=Foo
[Service]
RestartSec=22hr

