# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_4hour.target
[Unit]
Description=Foo
[Service]
WatchdogSec=4hour

