# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartupCPUWeight_10000.target
[Unit]
Description=Foo
[Service]
StartupCPUWeight=10000

