# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGPROF.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGPROF

