# systemdlint/tests/auto/InvalidValue/2.31/good/Service_MemoryDenyWriteExecute_1.target
[Unit]
Description=Foo
[Service]
MemoryDenyWriteExecute=1

