# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_8minute.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=8minute

