# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGURG.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGURG

