# systemdlint/tests/auto/InvalidValue/2.31/good/Service_ReadOnlyPaths_abc.target
[Unit]
Description=Foo
[Service]
ReadOnlyPaths=abc

