# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_6days.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=6days

