# systemdlint/tests/auto/InvalidValue/2.31/good/Service_MountAPIVFS_on.target
[Unit]
Description=Foo
[Service]
MountAPIVFS=on

