# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_16week.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=16week

