# systemdlint/tests/auto/InvalidValue/2.31/good/Service_ProtectSystem_on.target
[Unit]
Description=Foo
[Service]
ProtectSystem=on

