# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EWOULDBLOCK.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EWOULDBLOCK

