# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestrictSUIDSGID_yes.target
[Unit]
Description=Foo
[Service]
RestrictSUIDSGID=yes

