# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SyslogFacility_uucp.target
[Unit]
Description=Foo
[Service]
SyslogFacility=uucp

