# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_0days.target
[Unit]
Description=Foo
[Service]
WatchdogSec=0days

