# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGUNUSED.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGUNUSED

