# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGKILL.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGKILL

