# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_infinity.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=infinity

