# systemdlint/tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGIOT.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGIOT

