# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGIOT.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGIOT

