# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGTERM.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGTERM

