# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_0sec.target
[Unit]
Description=Foo
[Service]
TimeoutSec=0sec

