# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SyslogFacility_auth.target
[Unit]
Description=Foo
[Service]
SyslogFacility=auth

