# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_24years.target
[Unit]
Description=Foo
[Service]
TimeoutSec=24years

