# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SmackProcessLabel_abc.target
[Unit]
Description=Foo
[Service]
SmackProcessLabel=abc

