# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGPROF.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGPROF

