# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGTERM.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGTERM

