# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_1.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=1

