# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_20usec.target
[Unit]
Description=Foo
[Service]
WatchdogSec=20usec

