# systemdlint/tests/auto/InvalidValue/2.31/good/Service_MemoryLow_9999999999999G.target
[Unit]
Description=Foo
[Service]
MemoryLow=9999999999999G

