# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_16week.target
[Unit]
Description=Foo
[Service]
RestartSec=16week

