# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_6h.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=6h

