# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGKILL.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGKILL

