# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ETOOMANYREFS.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ETOOMANYREFS

