# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitAS_9999999999999.target
[Unit]
Description=Foo
[Service]
LimitAS=9999999999999

