# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_5hour.target
[Unit]
Description=Foo
[Service]
RestartSec=5hour

