# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGWINCH.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGWINCH

