# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_26months.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=26months

