# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_27y.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=27y

