# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGUNUSED.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGUNUSED

