# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartLimitInterval_infinity.target
[Unit]
Description=Foo
[Service]
StartLimitInterval=infinity

