# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_11m.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=11m

