# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EUCLEAN.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EUCLEAN

