# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_SIGSTOP.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=SIGSTOP

