# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_12second.target
[Unit]
Description=Foo
[Service]
TimeoutSec=12second

