# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGIO.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGIO

