# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_0days.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=0days

