# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_15min.target
[Unit]
Description=Foo
[Service]
WatchdogSec=15min

