# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_2us.target
[Unit]
Description=Foo
[Service]
TimeoutSec=2us

