# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallArchitectures_abc.target
[Unit]
Description=Foo
[Service]
SystemCallArchitectures=abc

