# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_8weeks.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=8weeks

