# systemdlint/tests/auto/InvalidValue/2.31/good/Service_User_iamauser.target
[Unit]
Description=Foo
[Service]
User=iamauser

