# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SyslogFacility_cron.target
[Unit]
Description=Foo
[Service]
SyslogFacility=cron

