# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGKILL.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGKILL

