# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ESPIPE.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ESPIPE

