# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGQUIT.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGQUIT

