# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_26year.target
[Unit]
Description=Foo
[Service]
WatchdogSec=26year

