# systemdlint/tests/auto/InvalidValue/2.31/good/Service_MemoryDenyWriteExecute_yes.target
[Unit]
Description=Foo
[Service]
MemoryDenyWriteExecute=yes

