# systemdlint/tests/auto/InvalidValue/2.31/good/Service_UMask_0.target
[Unit]
Description=Foo
[Service]
UMask=0

