# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartLimitInterval_22hr.target
[Unit]
Description=Foo
[Service]
StartLimitInterval=22hr

