# systemdlint/tests/auto/InvalidValue/2.31/good/Service_CPUAffinity_numa.target
[Unit]
Description=Foo
[Service]
CPUAffinity=numa

