# systemdlint/tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGILL.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGILL

