# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_8minute.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=8minute

