# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_24w.target
[Unit]
Description=Foo
[Service]
TimeoutSec=24w

