# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_25years.target
[Unit]
Description=Foo
[Service]
RestartSec=25years

