# systemdlint/tests/auto/InvalidValue/2.31/good/Service_PrivateDevices_on.target
[Unit]
Description=Foo
[Service]
PrivateDevices=on

