# systemdlint/tests/auto/InvalidValue/2.31/good/Service_Restart_always.target
[Unit]
Description=Foo
[Service]
Restart=always

