# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_6days.target
[Unit]
Description=Foo
[Service]
WatchdogSec=6days

