# systemdlint/tests/auto/InvalidValue/2.31/good/Service_GuessMainPID_true.target
[Unit]
Description=Foo
[Service]
GuessMainPID=true

