# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitMSGQUEUE_0.target
[Unit]
Description=Foo
[Service]
LimitMSGQUEUE=0

