# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGQUIT.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGQUIT

