# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGTTIN.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGTTIN

