# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_0sec.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=0sec

