# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_10day.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=10day

