# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_28h.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=28h

