# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_9min.target
[Unit]
Description=Foo
[Service]
WatchdogSec=9min

