# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_8.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=8

