# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGTRAP.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGTRAP

