# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGHUP.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGHUP

