# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EBADMSG.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EBADMSG

