# systemdlint/tests/auto/InvalidValue/2.31/good/Service_BlockIOAccounting_false.target
[Unit]
Description=Foo
[Service]
BlockIOAccounting=false

