# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGTERM.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGTERM

