# systemdlint/tests/auto/InvalidValue/2.31/good/Service_CPUAffinity_abc.target
[Unit]
Description=Foo
[Service]
CPUAffinity=abc

