# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_infinity.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=infinity

