# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGILL.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGILL

