# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_20w.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=20w

