# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_28h.target
[Unit]
Description=Foo
[Service]
RestartSec=28h

