# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_8weeks.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=8weeks

