# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_54.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=54

