# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_9.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=9

