# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_11y.target
[Unit]
Description=Foo
[Service]
WatchdogSec=11y

