# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TTYVHangup_true.target
[Unit]
Description=Foo
[Service]
TTYVHangup=true

