# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EUSERS.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EUSERS

