# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_19s.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=19s

