# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartLimitIntervalSec_26year.target
[Unit]
Description=Foo
[Service]
StartLimitIntervalSec=26year

