# systemdlint/tests/auto/InvalidValue/2.31/good/Service_MemoryAccounting_0.target
[Unit]
Description=Foo
[Service]
MemoryAccounting=0

