# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TTYVHangup_off.target
[Unit]
Description=Foo
[Service]
TTYVHangup=off

