# systemdlint/tests/auto/InvalidValue/2.31/good/Service_MemoryHigh_9999999999999K.target
[Unit]
Description=Foo
[Service]
MemoryHigh=9999999999999K

