# systemdlint/tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGXFSZ.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGXFSZ

