# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_2.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=2

