# systemdlint/tests/auto/InvalidValue/2.31/good/Service_MemorySwapMax_0B.target
[Unit]
Description=Foo
[Service]
MemorySwapMax=0B

