# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_11y.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=11y

