# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_25years.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=25years

