# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_3hours.target
[Unit]
Description=Foo
[Service]
RestartSec=3hours

