# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_20w.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=20w

