# systemdlint/tests/auto/InvalidValue/2.31/good/Service_CPUSchedulingPolicy_fifo.target
[Unit]
Description=Foo
[Service]
CPUSchedulingPolicy=fifo

