# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestrictRealtime_0.target
[Unit]
Description=Foo
[Service]
RestrictRealtime=0

