# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGPROF.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGPROF

