# systemdlint/tests/auto/InvalidValue/2.31/good/Service_CPUSchedulingPolicy_other.target
[Unit]
Description=Foo
[Service]
CPUSchedulingPolicy=other

