# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_27msec.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=27msec

