# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SupplementaryGroups_metoo.target
[Unit]
Description=Foo
[Service]
SupplementaryGroups=metoo

