# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_17second.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=17second

