# systemdlint/tests/auto/InvalidValue/2.31/good/Service_UtmpIdentifier_abc.target
[Unit]
Description=Foo
[Service]
UtmpIdentifier=abc

