# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_4hour.target
[Unit]
Description=Foo
[Service]
TimeoutSec=4hour

