# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitFSIZE_9999999999999.target
[Unit]
Description=Foo
[Service]
LimitFSIZE=9999999999999

