# systemdlint/tests/auto/InvalidValue/2.31/good/Service_Delegate_pids.target
[Unit]
Description=Foo
[Service]
Delegate=pids

