# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGTSTP.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGTSTP

