# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RootDirectory_--bin.target
[Unit]
Description=Foo
[Service]
RootDirectory=/bin

