# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartLimitIntervalSec_19s.target
[Unit]
Description=Foo
[Service]
StartLimitIntervalSec=19s

