# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGSYS.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGSYS

