# systemdlint/tests/auto/InvalidValue/2.31/good/Service_PrivateTmp_on.target
[Unit]
Description=Foo
[Service]
PrivateTmp=on

