# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EXFULL.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EXFULL

