# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_6h.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=6h

