# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_9month.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=9month

