# systemdlint/tests/auto/InvalidValue/2.31/good/Service_ProtectHome_true.target
[Unit]
Description=Foo
[Service]
ProtectHome=true

