# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StandardOutput_syslog.target
[Unit]
Description=Foo
[Service]
StandardOutput=syslog

