# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_7seconds.target
[Unit]
Description=Foo
[Service]
WatchdogSec=7seconds

