# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FileDescriptorStoreMax_0.target
[Unit]
Description=Foo
[Service]
FileDescriptorStoreMax=0

