# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_7minutes.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=7minutes

