# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_24w.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=24w

