# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGSTKFLT.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGSTKFLT

