# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGFPE.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGFPE

