# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitMSGQUEUE_9999999999999.target
[Unit]
Description=Foo
[Service]
LimitMSGQUEUE=9999999999999

