# systemdlint/tests/auto/InvalidValue/2.31/good/Service_IOWriteIOPSMax_abc.target
[Unit]
Description=Foo
[Service]
IOWriteIOPSMax=abc

