# systemdlint/tests/auto/InvalidValue/2.31/good/Service_CPUSchedulingPolicy_idle.target
[Unit]
Description=Foo
[Service]
CPUSchedulingPolicy=idle

