# systemdlint/tests/auto/InvalidValue/2.31/good/Service_InaccessableDirectories_abc.target
[Unit]
Description=Foo
[Service]
InaccessableDirectories=abc

