# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_1.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=1

