# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartLimitInterval_28h.target
[Unit]
Description=Foo
[Service]
StartLimitInterval=28h

