# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EIDRM.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EIDRM

