# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_22weeks.target
[Unit]
Description=Foo
[Service]
RestartSec=22weeks

