# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TasksMax_infinity.target
[Unit]
Description=Foo
[Service]
TasksMax=infinity

