# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_55.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=55

