# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_32.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=32

