# systemdlint/tests/auto/InvalidValue/2.31/good/Service_PIDFile_--usr--sbin.target
[Unit]
Description=Foo
[Service]
PIDFile=/usr/sbin

