# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SyslogFacility_authpriv.target
[Unit]
Description=Foo
[Service]
SyslogFacility=authpriv

