# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGURG.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGURG

