# systemdlint/tests/auto/InvalidValue/2.31/good/Service_MemoryMax_0K.target
[Unit]
Description=Foo
[Service]
MemoryMax=0K

