# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGTERM.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGTERM

