# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RemoveIPC_off.target
[Unit]
Description=Foo
[Service]
RemoveIPC=off

