# systemdlint/tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGUSR1.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGUSR1

