# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_24w.target
[Unit]
Description=Foo
[Service]
RestartSec=24w

