# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_8weeks.target
[Unit]
Description=Foo
[Service]
RestartSec=8weeks

