# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGLOST.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGLOST

