# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_0sec.target
[Unit]
Description=Foo
[Service]
RestartSec=0sec

