# systemdlint/tests/auto/InvalidValue/2.31/good/Service_ProtectSystem_strict.target
[Unit]
Description=Foo
[Service]
ProtectSystem=strict

