# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_15ms.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=15ms

