# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_DATAERR.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=DATAERR

