# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGABRT.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGABRT

