# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_0.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=0

