# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_15min.target
[Unit]
Description=Foo
[Service]
TimeoutSec=15min

