# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EMULTIHOP.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EMULTIHOP

