# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartLimitIntervalSec_6h.target
[Unit]
Description=Foo
[Service]
StartLimitIntervalSec=6h

