# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_12month.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=12month

