# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGXCPU.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGXCPU

