# systemdlint/tests/auto/InvalidValue/2.31/good/Service_IgnoreSIGPIPE_false.target
[Unit]
Description=Foo
[Service]
IgnoreSIGPIPE=false

