# systemdlint/tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGBUS.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGBUS

