# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_14M.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=14M

