# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_19s.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=19s

