# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitFSIZE_infinity.target
[Unit]
Description=Foo
[Service]
LimitFSIZE=infinity

