# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_28h.target
[Unit]
Description=Foo
[Service]
WatchdogSec=28h

