# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartupCPUShares_262144.target
[Unit]
Description=Foo
[Service]
StartupCPUShares=262144

