# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_6.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=6

