# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_6h.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=6h

