# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_1%C2%B5s.target
[Unit]
Description=Foo
[Service]
WatchdogSec=1µs

