# systemdlint/tests/auto/InvalidValue/2.31/good/Service_PrivateNetwork_1.target
[Unit]
Description=Foo
[Service]
PrivateNetwork=1

