# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EMFILE.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EMFILE

