# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_25years.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=25years

