# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_27y.target
[Unit]
Description=Foo
[Service]
WatchdogSec=27y

