# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_23m.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=23m

