# systemdlint/tests/auto/InvalidValue/2.31/good/Service_DeviceAllow_--dev--bar.target
[Unit]
Description=Foo
[Service]
DeviceAllow=/dev/bar

