# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartLimitIntervalSec_7minutes.target
[Unit]
Description=Foo
[Service]
StartLimitIntervalSec=7minutes

