# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_5hour.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=5hour

