# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StandardOutput_inherit.target
[Unit]
Description=Foo
[Service]
StandardOutput=inherit

