# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_8minute.target
[Unit]
Description=Foo
[Service]
WatchdogSec=8minute

