# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EDEADLOCK.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EDEADLOCK

