# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitRTTIME_0.target
[Unit]
Description=Foo
[Service]
LimitRTTIME=0

