# systemdlint/tests/auto/InvalidValue/2.31/good/Service_PrivateDevices_off.target
[Unit]
Description=Foo
[Service]
PrivateDevices=off

