# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGTTOU.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGTTOU

