# systemdlint/tests/auto/InvalidValue/2.31/good/Service_MemoryLow_0T.target
[Unit]
Description=Foo
[Service]
MemoryLow=0T

