# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_12second.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=12second

