# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGCLD.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGCLD

