# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TasksAccounting_on.target
[Unit]
Description=Foo
[Service]
TasksAccounting=on

