# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RootDirectoryStartOnly_off.target
[Unit]
Description=Foo
[Service]
RootDirectoryStartOnly=off

