# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EISNAM.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EISNAM

