# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_26year.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=26year

