# systemdlint/tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGABRT.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGABRT

