# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_18year.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=18year

