# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ETIME.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ETIME

