# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_4hour.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=4hour

