# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EBADE.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EBADE

