# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SyslogLevel_emerg.target
[Unit]
Description=Foo
[Service]
SyslogLevel=emerg

