# systemdlint/tests/auto/InvalidValue/2.31/good/Service_ProtectSystem_no.target
[Unit]
Description=Foo
[Service]
ProtectSystem=no

