# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WorkingDirectory_--bin.target
[Unit]
Description=Foo
[Service]
WorkingDirectory=/bin

