# systemdlint/tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGPOLL.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGPOLL

