# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_SIGLOST.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=SIGLOST

