# systemdlint/tests/auto/InvalidValue/2.31/good/Service_ReadWritePaths_abc.target
[Unit]
Description=Foo
[Service]
ReadWritePaths=abc

