# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_24.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=24

