# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_0sec.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=0sec

