# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_22weeks.target
[Unit]
Description=Foo
[Service]
WatchdogSec=22weeks

