# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_4ms.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=4ms

