# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGPOLL.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGPOLL

