# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_12second.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=12second

