# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_15min.target
[Unit]
Description=Foo
[Service]
RestartSec=15min

