# systemdlint/tests/auto/InvalidValue/2.31/good/Service_GuessMainPID_0.target
[Unit]
Description=Foo
[Service]
GuessMainPID=0

