# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_17minutes.target
[Unit]
Description=Foo
[Service]
TimeoutSec=17minutes

