# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TTYVHangup_1.target
[Unit]
Description=Foo
[Service]
TTYVHangup=1

