# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_7seconds.target
[Unit]
Description=Foo
[Service]
TimeoutSec=7seconds

