# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_19d.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=19d

