# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_41.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=41

