# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_22hr.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=22hr

