# systemdlint/tests/auto/InvalidValue/2.31/good/Service_Group_iamagroup.target
[Unit]
Description=Foo
[Service]
Group=iamagroup

