# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_9999999999999.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=9999999999999

