# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_11y.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=11y

