# systemdlint/tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGQUIT.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGQUIT

