# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_25s.target
[Unit]
Description=Foo
[Service]
TimeoutSec=25s

