# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGTTIN.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGTTIN

