# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_16seconds.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=16seconds

