# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_10day.target
[Unit]
Description=Foo
[Service]
WatchdogSec=10day

