# systemdlint/tests/auto/InvalidValue/2.31/good/Service_NonBlocking_no.target
[Unit]
Description=Foo
[Service]
NonBlocking=no

