# systemdlint/tests/auto/InvalidValue/2.31/good/Service_MemoryLimit_infinity.target
[Unit]
Description=Foo
[Service]
MemoryLimit=infinity

