# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ENFILE.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ENFILE

