# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_6days.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=6days

