# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EBADR.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EBADR

