# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EXDEV.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EXDEV

