# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGSYS.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGSYS

