# systemdlint/tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGINFO.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGINFO

