# systemdlint/tests/auto/InvalidValue/2.31/good/Service_CPUSchedulingResetOnFork_1.target
[Unit]
Description=Foo
[Service]
CPUSchedulingResetOnFork=1

