# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RootDirectoryStartOnly_0.target
[Unit]
Description=Foo
[Service]
RootDirectoryStartOnly=0

