# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_19d.target
[Unit]
Description=Foo
[Service]
TimeoutSec=19d

