# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGXFSZ.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGXFSZ

