# systemdlint/tests/auto/InvalidValue/2.31/good/Service_Restart_on-watchdog.target
[Unit]
Description=Foo
[Service]
Restart=on-watchdog

