# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitNICE_9999999999999.target
[Unit]
Description=Foo
[Service]
LimitNICE=9999999999999

