# systemdlint/tests/auto/InvalidValue/2.31/good/Service_IgnoreSIGPIPE_true.target
[Unit]
Description=Foo
[Service]
IgnoreSIGPIPE=true

