# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGSEGV.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGSEGV

