# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_16seconds.target
[Unit]
Description=Foo
[Service]
WatchdogSec=16seconds

