# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EIO.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EIO

