# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SyslogLevel_info.target
[Unit]
Description=Foo
[Service]
SyslogLevel=info

