# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StandardError_abc.target
[Unit]
Description=Foo
[Service]
StandardError=abc

