# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_13minute.target
[Unit]
Description=Foo
[Service]
WatchdogSec=13minute

