# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGXCPU.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGXCPU

