# systemdlint/tests/auto/InvalidValue/2.31/good/Service_BlockIOWriteBandwidthMax_--dev--bar%209999999999999G.target
[Unit]
Description=Foo
[Service]
BlockIOWriteBandwidthMax=/dev/bar 9999999999999G

