# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGXCPU.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGXCPU

