# systemdlint/tests/auto/InvalidValue/2.31/good/Service_ExecStartPost_--bin--sh.target
[Unit]
Description=Foo
[Service]
ExecStartPost=/bin/sh

