# systemdlint/tests/auto/InvalidValue/2.31/good/Service_MountFlags_shared.target
[Unit]
Description=Foo
[Service]
MountFlags=shared

