# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGINFO.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGINFO

