# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_23m.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=23m

