# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_28%C2%B5s.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=28µs

