# systemdlint/tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGTSTP.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGTSTP

