# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TasksAccounting_no.target
[Unit]
Description=Foo
[Service]
TasksAccounting=no

