# systemdlint/tests/auto/InvalidValue/2.31/good/Service_Nice_-20.target
[Unit]
Description=Foo
[Service]
Nice=-20

