# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EISDIR.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EISDIR

