# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WorkingDirectory_--usr--sbin.target
[Unit]
Description=Foo
[Service]
WorkingDirectory=/usr/sbin

