# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RootDirectoryStartOnly_false.target
[Unit]
Description=Foo
[Service]
RootDirectoryStartOnly=false

