# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGBUS.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGBUS

