# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_1day.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=1day

