# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_24w.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=24w

