# systemdlint/tests/auto/InvalidValue/2.31/good/Service_PAMName_abc.target
[Unit]
Description=Foo
[Service]
PAMName=abc

