# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_39.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=39

