# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_13minute.target
[Unit]
Description=Foo
[Service]
TimeoutSec=13minute

