# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGUSR1.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGUSR1

