# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_27y.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=27y

