# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGKILL.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGKILL

