# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_SIGXCPU.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=SIGXCPU

