# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_21usec.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=21usec

