# systemdlint/tests/auto/InvalidValue/2.31/good/Service_IODeviceWeight_abc.target
[Unit]
Description=Foo
[Service]
IODeviceWeight=abc

