# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_0.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=0

