# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGINT.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGINT

