# systemdlint/tests/auto/InvalidValue/2.31/good/Service_NotifyAccess_none.target
[Unit]
Description=Foo
[Service]
NotifyAccess=none

