# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestrictSUIDSGID_off.target
[Unit]
Description=Foo
[Service]
RestrictSUIDSGID=off

