# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ENOBUFS.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ENOBUFS

