# systemdlint/tests/auto/InvalidValue/2.31/good/Service_NonBlocking_yes.target
[Unit]
Description=Foo
[Service]
NonBlocking=yes

