# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestrictSUIDSGID_no.target
[Unit]
Description=Foo
[Service]
RestrictSUIDSGID=no

