# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_11m.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=11m

