# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_25years.target
[Unit]
Description=Foo
[Service]
WatchdogSec=25years

