# systemdlint/tests/auto/InvalidValue/2.31/good/Service_IgnoreSIGPIPE_on.target
[Unit]
Description=Foo
[Service]
IgnoreSIGPIPE=on

