# systemdlint/tests/auto/InvalidValue/2.31/good/Service_BlockStartupIOWeight_10000.target
[Unit]
Description=Foo
[Service]
BlockStartupIOWeight=10000

