# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_6days.target
[Unit]
Description=Foo
[Service]
TimeoutSec=6days

