# systemdlint/tests/auto/InvalidValue/2.31/good/Service_MemoryDenyWriteExecute_false.target
[Unit]
Description=Foo
[Service]
MemoryDenyWriteExecute=false

