# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_12month.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=12month

