# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitNPROC_0.target
[Unit]
Description=Foo
[Service]
LimitNPROC=0

