# systemdlint/tests/auto/InvalidValue/2.31/good/Service_CPUSchedulingPolicy_rr.target
[Unit]
Description=Foo
[Service]
CPUSchedulingPolicy=rr

