# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StandardOutput_abc.target
[Unit]
Description=Foo
[Service]
StandardOutput=abc

