# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGFPE.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGFPE

