# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGUSR2.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGUSR2

