# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TTYVTDisallocate_no.target
[Unit]
Description=Foo
[Service]
TTYVTDisallocate=no

