# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartLimitIntervalSec_9999999999999.target
[Unit]
Description=Foo
[Service]
StartLimitIntervalSec=9999999999999

