# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_2us.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=2us

