# systemdlint/tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGEMT.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGEMT

