# systemdlint/tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGSTKFLT.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGSTKFLT

