# systemdlint/tests/auto/InvalidValue/2.31/good/Service_CPUSchedulingPriority_99.target
[Unit]
Description=Foo
[Service]
CPUSchedulingPriority=99

