# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_16week.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=16week

