# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_10day.target
[Unit]
Description=Foo
[Service]
RestartSec=10day

