# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartLimitBurst_9999999999999.target
[Unit]
Description=Foo
[Service]
StartLimitBurst=9999999999999

