# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_28h.target
[Unit]
Description=Foo
[Service]
TimeoutSec=28h

