# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartLimitInterval_4ms.target
[Unit]
Description=Foo
[Service]
StartLimitInterval=4ms

