# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_6days.target
[Unit]
Description=Foo
[Service]
RestartSec=6days

