# systemdlint/tests/auto/InvalidValue/2.31/good/Service_PrivateTmp_off.target
[Unit]
Description=Foo
[Service]
PrivateTmp=off

