# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ELIBMAX.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ELIBMAX

