# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_14M.target
[Unit]
Description=Foo
[Service]
TimeoutSec=14M

