# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_23m.target
[Unit]
Description=Foo
[Service]
TimeoutSec=23m

