# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_27msec.target
[Unit]
Description=Foo
[Service]
RestartSec=27msec

