# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_8minute.target
[Unit]
Description=Foo
[Service]
TimeoutSec=8minute

