# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TTYVHangup_yes.target
[Unit]
Description=Foo
[Service]
TTYVHangup=yes

