# systemdlint/tests/auto/InvalidValue/2.31/good/Service_PrivateDevices_1.target
[Unit]
Description=Foo
[Service]
PrivateDevices=1

