# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_1%C2%B5s.target
[Unit]
Description=Foo
[Service]
RestartSec=1µs

