# systemdlint/tests/auto/InvalidValue/2.31/good/Service_IgnoreSIGPIPE_off.target
[Unit]
Description=Foo
[Service]
IgnoreSIGPIPE=off

