# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_27.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=27

