# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TTYReset_true.target
[Unit]
Description=Foo
[Service]
TTYReset=true

