# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGINFO.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGINFO

