# systemdlint/tests/auto/InvalidValue/2.31/good/Service_MemoryAccounting_no.target
[Unit]
Description=Foo
[Service]
MemoryAccounting=no

