# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitNICE_infinity.target
[Unit]
Description=Foo
[Service]
LimitNICE=infinity

