# systemdlint/tests/auto/InvalidValue/2.31/good/Service_ProtectHome_yes.target
[Unit]
Description=Foo
[Service]
ProtectHome=yes

