# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_DATAERR.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=DATAERR

