# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_26year.target
[Unit]
Description=Foo
[Service]
RestartSec=26year

