# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitRSS_9999999999999.target
[Unit]
Description=Foo
[Service]
LimitRSS=9999999999999

