# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_28h.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=28h

