# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_2us.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=2us

