# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_5hour.target
[Unit]
Description=Foo
[Service]
WatchdogSec=5hour

