# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ENOLINK.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ENOLINK

