# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SendSIGKILL_false.target
[Unit]
Description=Foo
[Service]
SendSIGKILL=false

