# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_16seconds.target
[Unit]
Description=Foo
[Service]
TimeoutSec=16seconds

