# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_3.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=3

