# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RemoveIPC_on.target
[Unit]
Description=Foo
[Service]
RemoveIPC=on

