# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGBUS.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGBUS

