# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_2us.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=2us

