# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGSTOP.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGSTOP

