# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_18sec.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=18sec

