# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_SIGPIPE.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=SIGPIPE

