# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TTYVTDisallocate_false.target
[Unit]
Description=Foo
[Service]
TTYVTDisallocate=false

