# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_19d.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=19d

