# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartLimitIntervalSec_24w.target
[Unit]
Description=Foo
[Service]
StartLimitIntervalSec=24w

