# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EPIPE.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EPIPE

