# systemdlint/tests/auto/InvalidValue/2.31/good/Service_LimitNOFILE_9999999999999.target
[Unit]
Description=Foo
[Service]
LimitNOFILE=9999999999999

