# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_8weeks.target
[Unit]
Description=Foo
[Service]
TimeoutSec=8weeks

