# systemdlint/tests/auto/InvalidValue/2.31/good/Service_AmbientCapabilities_abc.target
[Unit]
Description=Foo
[Service]
AmbientCapabilities=abc

