# systemdlint/tests/auto/InvalidValue/2.31/good/Service_Delegate_cpu.target
[Unit]
Description=Foo
[Service]
Delegate=cpu

