# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_11y.target
[Unit]
Description=Foo
[Service]
RestartSec=11y

