# systemdlint/tests/auto/InvalidValue/2.31/good/Service_StartLimitInterval_27y.target
[Unit]
Description=Foo
[Service]
StartLimitInterval=27y

