# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_18year.target
[Unit]
Description=Foo
[Service]
TimeoutSec=18year

