# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_27msec.target
[Unit]
Description=Foo
[Service]
TimeoutSec=27msec

