# systemdlint/tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGBUS.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGBUS

