# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_SIGVTALRM.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=SIGVTALRM

