# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_27y.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=27y

