# systemdlint/tests/auto/InvalidValue/2.31/good/Service_PrivateDevices_no.target
[Unit]
Description=Foo
[Service]
PrivateDevices=no

