# systemdlint/tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_7minutes.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=7minutes

