# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SendSIGHUP_no.target
[Unit]
Description=Foo
[Service]
SendSIGHUP=no

