# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RestartSec_23week.target
[Unit]
Description=Foo
[Service]
RestartSec=23week

