# systemdlint/tests/auto/InvalidValue/2.31/good/Service_Personality_s390.target
[Unit]
Description=Foo
[Service]
Personality=s390

