# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_46.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=46

