# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RootDirectoryStartOnly_on.target
[Unit]
Description=Foo
[Service]
RootDirectoryStartOnly=on

