# systemdlint/tests/auto/InvalidValue/2.31/good/Service_ExecStart_--bin--sh.target
[Unit]
Description=Foo
[Service]
ExecStart=/bin/sh

