# systemdlint/tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_24w.target
[Unit]
Description=Foo
[Service]
WatchdogSec=24w

