# systemdlint/tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_26year.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=26year

