# systemdlint/tests/auto/InvalidValue/2.31/good/Service_SyslogLevel_debug.target
[Unit]
Description=Foo
[Service]
SyslogLevel=debug

