# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TTYReset_on.target
[Unit]
Description=Foo
[Socket]
TTYReset=on

