# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGCHLD.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGCHLD

