# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SyslogFacility_uucp.target
[Unit]
Description=Foo
[Socket]
SyslogFacility=uucp

