# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TimeoutSec_21usec.target
[Unit]
Description=Foo
[Socket]
TimeoutSec=21usec

