# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KeepAlive_1.target
[Unit]
Description=Foo
[Socket]
KeepAlive=1

