# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SELinuxContextFromNet_off.target
[Unit]
Description=Foo
[Socket]
SELinuxContextFromNet=off

