# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_RestrictRealtime_no.target
[Unit]
Description=Foo
[Socket]
RestrictRealtime=no

