# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SmackLabelIPIn_abc.target
[Unit]
Description=Foo
[Socket]
SmackLabelIPIn=abc

