# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EHWPOISON.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EHWPOISON

