# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_DisableControllers_io.target
[Unit]
Description=Foo
[Socket]
DisableControllers=io

