# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGURG.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGURG

