# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_BlockIOAccounting_off.target
[Unit]
Description=Foo
[Socket]
BlockIOAccounting=off

