# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitIntervalSec_8weeks.target
[Unit]
Description=Foo
[Socket]
TriggerLimitIntervalSec=8weeks

