# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_MemoryAccounting_no.target
[Unit]
Description=Foo
[Socket]
MemoryAccounting=no

