# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_34.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=34

