# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_ProtectHome_on.target
[Unit]
Description=Foo
[Socket]
ProtectHome=on

