# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitMEMLOCK_0.target
[Unit]
Description=Foo
[Socket]
LimitMEMLOCK=0

