# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TimeoutSec_6days.target
[Unit]
Description=Foo
[Socket]
TimeoutSec=6days

