# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_18year.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=18year

