# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KeepAliveTimeSec_0.target
[Unit]
Description=Foo
[Socket]
KeepAliveTimeSec=0

