# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_StandardError_tty.target
[Unit]
Description=Foo
[Socket]
StandardError=tty

