# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_AmbientCapabilities_abc.target
[Unit]
Description=Foo
[Socket]
AmbientCapabilities=abc

