# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TimeoutSec_19d.target
[Unit]
Description=Foo
[Socket]
TimeoutSec=19d

