# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitDATA_infinity.target
[Unit]
Description=Foo
[Socket]
LimitDATA=infinity

