# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_ReceiveBuffer_9999999999999M.target
[Unit]
Description=Foo
[Socket]
ReceiveBuffer=9999999999999M

