# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TimeoutSec_24w.target
[Unit]
Description=Foo
[Socket]
TimeoutSec=24w

