# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_3.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=3

