# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_PassCredentials_false.target
[Unit]
Description=Foo
[Socket]
PassCredentials=false

