# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_CPUSchedulingPriority_99.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingPriority=99

