# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitSIGPENDING_9999999999999.target
[Unit]
Description=Foo
[Socket]
LimitSIGPENDING=9999999999999

