# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitNOFILE_9999999999999.target
[Unit]
Description=Foo
[Socket]
LimitNOFILE=9999999999999

