# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_1day.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=1day

