# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SyslogLevelPrefix_off.target
[Unit]
Description=Foo
[Socket]
SyslogLevelPrefix=off

