# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_DeviceAllow_--dev--bar.target
[Unit]
Description=Foo
[Socket]
DeviceAllow=/dev/bar

