# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOTDIR.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOTDIR

