# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGBUS.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGBUS

