# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGTSTP.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGTSTP

