# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_PassSecurity_on.target
[Unit]
Description=Foo
[Socket]
PassSecurity=on

