# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_Mark_0.target
[Unit]
Description=Foo
[Socket]
Mark=0

