# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_IOWeight_1.target
[Unit]
Description=Foo
[Socket]
IOWeight=1

