# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EAGAIN.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EAGAIN

