# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGEMT.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGEMT

