# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EDESTADDRREQ.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EDESTADDRREQ

