# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_CPUSchedulingPolicy_idle.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingPolicy=idle

