# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_RestrictRealtime_true.target
[Unit]
Description=Foo
[Socket]
RestrictRealtime=true

