# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_44.target
[Unit]
Description=Foo
[Socket]
KillSignal=44

