# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SecureBits_keep-caps-locked.target
[Unit]
Description=Foo
[Socket]
SecureBits=keep-caps-locked

