# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_Accept_off.target
[Unit]
Description=Foo
[Socket]
Accept=off

