# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KeepAlive_true.target
[Unit]
Description=Foo
[Socket]
KeepAlive=true

