# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ESHUTDOWN.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ESHUTDOWN

