# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_MemoryLimit_0K.target
[Unit]
Description=Foo
[Socket]
MemoryLimit=0K

