# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TTYVHangup_off.target
[Unit]
Description=Foo
[Socket]
TTYVHangup=off

