# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGURG.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGURG

