# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_ProtectHome_read-only.target
[Unit]
Description=Foo
[Socket]
ProtectHome=read-only

