# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_ListenUSBFunction_abc.target
[Unit]
Description=Foo
[Socket]
ListenUSBFunction=abc

