# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_10day.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=10day

