# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TTYVTDisallocate_on.target
[Unit]
Description=Foo
[Socket]
TTYVTDisallocate=on

