# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ERANGE.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ERANGE

