# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EISCONN.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EISCONN

