# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SupplementaryGroups_iamagroup.target
[Unit]
Description=Foo
[Socket]
SupplementaryGroups=iamagroup

