# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_38.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=38

