# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_Transparent_off.target
[Unit]
Description=Foo
[Socket]
Transparent=off

