# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGWINCH.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGWINCH

