# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_IODeviceLatencyTargetSec_abc.target
[Unit]
Description=Foo
[Socket]
IODeviceLatencyTargetSec=abc

