# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TimerSlackNSec_0.target
[Unit]
Description=Foo
[Socket]
TimerSlackNSec=0

