# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitMEMLOCK_infinity.target
[Unit]
Description=Foo
[Socket]
LimitMEMLOCK=infinity

