# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_MemoryLimit_infinity.target
[Unit]
Description=Foo
[Socket]
MemoryLimit=infinity

