# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_StandardError_syslog.target
[Unit]
Description=Foo
[Socket]
StandardError=syslog

