# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EBADFD.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EBADFD

