# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_2us.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=2us

