# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_StandardError_null.target
[Unit]
Description=Foo
[Socket]
StandardError=null

