# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_PassCredentials_on.target
[Unit]
Description=Foo
[Socket]
PassCredentials=on

