# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitFSIZE_0.target
[Unit]
Description=Foo
[Socket]
LimitFSIZE=0

