# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_BlockIOReadBandwidthMax_--dev--foo%209999999999999B.target
[Unit]
Description=Foo
[Socket]
BlockIOReadBandwidthMax=/dev/foo 9999999999999B

