# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TTYVTDisallocate_yes.target
[Unit]
Description=Foo
[Socket]
TTYVTDisallocate=yes

