# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGLOST.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGLOST

