# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_PrivateTmp_on.target
[Unit]
Description=Foo
[Socket]
PrivateTmp=on

