# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SmackLabel_abc.target
[Unit]
Description=Foo
[Socket]
SmackLabel=abc

