# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TTYVHangup_1.target
[Unit]
Description=Foo
[Socket]
TTYVHangup=1

