# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGINT.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGINT

