# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_StandardOutput_kmsg%2Bconsole.target
[Unit]
Description=Foo
[Socket]
StandardOutput=kmsg+console

