# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGVTALRM.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGVTALRM

