# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGTTOU.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGTTOU

