# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_DATAERR.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=DATAERR

