# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_PassSecurity_no.target
[Unit]
Description=Foo
[Socket]
PassSecurity=no

