# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGUNUSED.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGUNUSED

