# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGTRAP.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGTRAP

