# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SocketProtocol_sctp.target
[Unit]
Description=Foo
[Socket]
SocketProtocol=sctp

