# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SELinuxContextFromNet_1.target
[Unit]
Description=Foo
[Socket]
SELinuxContextFromNet=1

