# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ETOOMANYREFS.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ETOOMANYREFS

