# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_0days.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=0days

