# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGABRT.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGABRT

