# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ETIME.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ETIME

