# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_BlockIOAccounting_yes.target
[Unit]
Description=Foo
[Socket]
BlockIOAccounting=yes

