# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_RestrictAddressFamilies_abc.target
[Unit]
Description=Foo
[Socket]
RestrictAddressFamilies=abc

