# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_RemoveIPC_off.target
[Unit]
Description=Foo
[Socket]
RemoveIPC=off

