# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGSTOP.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGSTOP

