# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGUNUSED.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGUNUSED

