# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGUNUSED.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGUNUSED

