# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGCONT.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGCONT

