# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_PassSecurity_0.target
[Unit]
Description=Foo
[Socket]
PassSecurity=0

