# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_Delegate_memory.target
[Unit]
Description=Foo
[Socket]
Delegate=memory

