# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_25.target
[Unit]
Description=Foo
[Socket]
KillSignal=25

