# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_8minute.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=8minute

