# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TasksMax_9999999999999.target
[Unit]
Description=Foo
[Socket]
TasksMax=9999999999999

