# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_Transparent_no.target
[Unit]
Description=Foo
[Socket]
Transparent=no

