# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENODATA.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENODATA

