# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGUSR2.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGUSR2

