# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TimeoutSec_15ms.target
[Unit]
Description=Foo
[Socket]
TimeoutSec=15ms

