# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_IOSchedulingPriority_0.target
[Unit]
Description=Foo
[Socket]
IOSchedulingPriority=0

