# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_Writable_off.target
[Unit]
Description=Foo
[Socket]
Writable=off

