# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillMode_control-group.target
[Unit]
Description=Foo
[Socket]
KillMode=control-group

