# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_DeferAcceptSec_9999999999999.target
[Unit]
Description=Foo
[Socket]
DeferAcceptSec=9999999999999

