# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SyslogLevelPrefix_1.target
[Unit]
Description=Foo
[Socket]
SyslogLevelPrefix=1

