# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_CPUSchedulingResetOnFork_on.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingResetOnFork=on

