# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_DATAERR.target
[Unit]
Description=Foo
[Socket]
KillSignal=DATAERR

