# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SyslogFacility_mail.target
[Unit]
Description=Foo
[Socket]
SyslogFacility=mail

