# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOMEM.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOMEM

