# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_CPUSchedulingResetOnFork_false.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingResetOnFork=false

