# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_PrivateNetwork_1.target
[Unit]
Description=Foo
[Socket]
PrivateNetwork=1

