# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_MemoryDenyWriteExecute_off.target
[Unit]
Description=Foo
[Socket]
MemoryDenyWriteExecute=off

