# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_IPTOS_0.target
[Unit]
Description=Foo
[Socket]
IPTOS=0

