# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EDOM.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EDOM

