# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SyslogFacility_news.target
[Unit]
Description=Foo
[Socket]
SyslogFacility=news

