# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitNPROC_0.target
[Unit]
Description=Foo
[Socket]
LimitNPROC=0

