# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_Group_metoo.target
[Unit]
Description=Foo
[Socket]
Group=metoo

