# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SyslogFacility_authpriv.target
[Unit]
Description=Foo
[Socket]
SyslogFacility=authpriv

