# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_29.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=29

