# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGKILL.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGKILL

