# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SocketMode_644.target
[Unit]
Description=Foo
[Socket]
SocketMode=644

