# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SendSIGHUP_1.target
[Unit]
Description=Foo
[Socket]
SendSIGHUP=1

