# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGCLD.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGCLD

