# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGFPE.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGFPE

