# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_4095.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=4095

