# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_BlockIOAccounting_true.target
[Unit]
Description=Foo
[Socket]
BlockIOAccounting=true

