# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_13minute.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=13minute

