# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_ListenNetlink_abc.target
[Unit]
Description=Foo
[Socket]
ListenNetlink=abc

