# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_IgnoreSIGPIPE_0.target
[Unit]
Description=Foo
[Socket]
IgnoreSIGPIPE=0

