# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOTEMPTY.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOTEMPTY

