# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitCPU_0.target
[Unit]
Description=Foo
[Socket]
LimitCPU=0

