# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGPWR.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGPWR

