# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENONET.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENONET

