# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SendSIGHUP_yes.target
[Unit]
Description=Foo
[Socket]
SendSIGHUP=yes

