# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGXFSZ.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGXFSZ

