# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ESTALE.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ESTALE

