# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitIntervalSec_21usec.target
[Unit]
Description=Foo
[Socket]
TriggerLimitIntervalSec=21usec

