# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_7minutes.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=7minutes

