# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ECONNRESET.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ECONNRESET

