# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_Broadcast_yes.target
[Unit]
Description=Foo
[Socket]
Broadcast=yes

