# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_RestrictSUIDSGID_1.target
[Unit]
Description=Foo
[Socket]
RestrictSUIDSGID=1

