# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_IgnoreSIGPIPE_1.target
[Unit]
Description=Foo
[Socket]
IgnoreSIGPIPE=1

