# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TimeoutSec_12second.target
[Unit]
Description=Foo
[Socket]
TimeoutSec=12second

