# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGTTIN.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGTTIN

