# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KeepAlive_on.target
[Unit]
Description=Foo
[Socket]
KeepAlive=on

