# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_Environment_a%3Db.target
[Unit]
Description=Foo
[Socket]
Environment=a=b

