# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGFPE.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGFPE

