# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_InaccessiblePaths_abc.target
[Unit]
Description=Foo
[Socket]
InaccessiblePaths=abc

