# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_IPTOS_reliability.target
[Unit]
Description=Foo
[Socket]
IPTOS=reliability

