# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EOWNERDEAD.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EOWNERDEAD

