# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_PrivateTmp_false.target
[Unit]
Description=Foo
[Socket]
PrivateTmp=false

