# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitIntervalSec_22hr.target
[Unit]
Description=Foo
[Socket]
TriggerLimitIntervalSec=22hr

