# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_CPUSchedulingResetOnFork_1.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingResetOnFork=1

