# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TTYVHangup_false.target
[Unit]
Description=Foo
[Socket]
TTYVHangup=false

