# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_ProtectSystem_false.target
[Unit]
Description=Foo
[Socket]
ProtectSystem=false

