# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_NoDelay_false.target
[Unit]
Description=Foo
[Socket]
NoDelay=false

