# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_RemoveOnStop_no.target
[Unit]
Description=Foo
[Socket]
RemoveOnStop=no

