# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitCORE_9999999999999.target
[Unit]
Description=Foo
[Socket]
LimitCORE=9999999999999

