# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_21usec.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=21usec

