# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EALREADY.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EALREADY

