# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_8.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=8

