# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_StandardError_socket.target
[Unit]
Description=Foo
[Socket]
StandardError=socket

