# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SupplementaryGroups_metoo.target
[Unit]
Description=Foo
[Socket]
SupplementaryGroups=metoo

