# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_Broadcast_no.target
[Unit]
Description=Foo
[Socket]
Broadcast=no

