# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_IgnoreSIGPIPE_off.target
[Unit]
Description=Foo
[Socket]
IgnoreSIGPIPE=off

