# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGWINCH.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGWINCH

