# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_ProtectHome_true.target
[Unit]
Description=Foo
[Socket]
ProtectHome=true

