# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_IgnoreSIGPIPE_false.target
[Unit]
Description=Foo
[Socket]
IgnoreSIGPIPE=false

