# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SendSIGKILL_off.target
[Unit]
Description=Foo
[Socket]
SendSIGKILL=off

