# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOPROTOOPT.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOPROTOOPT

