# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_24.target
[Unit]
Description=Foo
[Socket]
KillSignal=24

