# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_IgnoreSIGPIPE_yes.target
[Unit]
Description=Foo
[Socket]
IgnoreSIGPIPE=yes

