# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_MemoryLimit_9999999999999B.target
[Unit]
Description=Foo
[Socket]
MemoryLimit=9999999999999B

