# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOSYS.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOSYS

