# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOEXEC.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOEXEC

