# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EHOSTDOWN.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EHOSTDOWN

