# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGCHLD.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGCHLD

