# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FreeBind_0.target
[Unit]
Description=Foo
[Socket]
FreeBind=0

