# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGVTALRM.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGVTALRM

