# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_CPUSchedulingResetOnFork_true.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingResetOnFork=true

