# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOBUFS.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOBUFS

