# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_ReusePort_no.target
[Unit]
Description=Foo
[Socket]
ReusePort=no

