# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitRTTIME_infinity.target
[Unit]
Description=Foo
[Socket]
LimitRTTIME=infinity

