# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitLOCKS_0.target
[Unit]
Description=Foo
[Socket]
LimitLOCKS=0

