# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EISNAM.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EISNAM

