# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TimeoutSec_infinity.target
[Unit]
Description=Foo
[Socket]
TimeoutSec=infinity

