# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_EnvironmentFile_--dev--.target
[Unit]
Description=Foo
[Socket]
EnvironmentFile=/dev/

