# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_RestrictRealtime_0.target
[Unit]
Description=Foo
[Socket]
RestrictRealtime=0

