# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGSTKFLT.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGSTKFLT

