# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitIntervalSec_6days.target
[Unit]
Description=Foo
[Socket]
TriggerLimitIntervalSec=6days

