# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOSR.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOSR

