# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_DevicePolicy_strict.target
[Unit]
Description=Foo
[Socket]
DevicePolicy=strict

