# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_PassSecurity_1.target
[Unit]
Description=Foo
[Socket]
PassSecurity=1

