# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EUSERS.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EUSERS

