# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_32.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=32

