# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENETUNREACH.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENETUNREACH

