# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_IPTOS_low-delay.target
[Unit]
Description=Foo
[Socket]
IPTOS=low-delay

