# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EBADR.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EBADR

