# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitSTACK_9999999999999.target
[Unit]
Description=Foo
[Socket]
LimitSTACK=9999999999999

