# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ELIBMAX.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ELIBMAX

