# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TTYVTDisallocate_true.target
[Unit]
Description=Foo
[Socket]
TTYVTDisallocate=true

