# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TimeoutSec_8minute.target
[Unit]
Description=Foo
[Socket]
TimeoutSec=8minute

