# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_E2BIG.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=E2BIG

