# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EROFS.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EROFS

