# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_StandardOutput_inherit.target
[Unit]
Description=Foo
[Socket]
StandardOutput=inherit

