# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_IOSchedulingClass_idle.target
[Unit]
Description=Foo
[Socket]
IOSchedulingClass=idle

