# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGBUS.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGBUS

