# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_MountAPIVFS_false.target
[Unit]
Description=Foo
[Socket]
MountAPIVFS=false

