# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_RemoveIPC_yes.target
[Unit]
Description=Foo
[Socket]
RemoveIPC=yes

