# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_ProtectHome_0.target
[Unit]
Description=Foo
[Socket]
ProtectHome=0

