# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EWOULDBLOCK.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EWOULDBLOCK

