# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TimeoutSec_25s.target
[Unit]
Description=Foo
[Socket]
TimeoutSec=25s

