# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_ProtectSystem_on.target
[Unit]
Description=Foo
[Socket]
ProtectSystem=on

