# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_BlockIOReadBandwidthMax_--dev--bar%200G.target
[Unit]
Description=Foo
[Socket]
BlockIOReadBandwidthMax=/dev/bar 0G

