# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TTYReset_yes.target
[Unit]
Description=Foo
[Socket]
TTYReset=yes

