# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_MemoryHigh_0G.target
[Unit]
Description=Foo
[Socket]
MemoryHigh=0G

