# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SyslogFacility_local0.target
[Unit]
Description=Foo
[Socket]
SyslogFacility=local0

