# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_Personality_x86.target
[Unit]
Description=Foo
[Socket]
Personality=x86

