# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_StartupCPUShares_2.target
[Unit]
Description=Foo
[Socket]
StartupCPUShares=2

