# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KeepAliveTimeSec_9999999999999.target
[Unit]
Description=Foo
[Socket]
KeepAliveTimeSec=9999999999999

