# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_StandardError_inherit.target
[Unit]
Description=Foo
[Socket]
StandardError=inherit

