# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_NoDelay_on.target
[Unit]
Description=Foo
[Socket]
NoDelay=on

