# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_MemoryDenyWriteExecute_on.target
[Unit]
Description=Foo
[Socket]
MemoryDenyWriteExecute=on

