# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_ProtectSystem_0.target
[Unit]
Description=Foo
[Socket]
ProtectSystem=0

