# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_PrivateDevices_false.target
[Unit]
Description=Foo
[Socket]
PrivateDevices=false

