# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_StandardOutput_syslog.target
[Unit]
Description=Foo
[Socket]
StandardOutput=syslog

