# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_CPUSchedulingPolicy_fifo.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingPolicy=fifo

