# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KeepAliveProbes_9999999999999.target
[Unit]
Description=Foo
[Socket]
KeepAliveProbes=9999999999999

