# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SocketUser_iamauser.target
[Unit]
Description=Foo
[Socket]
SocketUser=iamauser

