# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_StandardOutput_syslog%2Bconsole.target
[Unit]
Description=Foo
[Socket]
StandardOutput=syslog+console

