# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_ProtectHome_1.target
[Unit]
Description=Foo
[Socket]
ProtectHome=1

