# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_ExecStopPre_--bin--sh.target
[Unit]
Description=Foo
[Socket]
ExecStopPre=/bin/sh

