# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_UtmpMode_login.target
[Unit]
Description=Foo
[Socket]
UtmpMode=login

