# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_3hours.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=3hours

