# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_RemoveOnStop_true.target
[Unit]
Description=Foo
[Socket]
RemoveOnStop=true

