# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_MessageQueueMessageSize_0.target
[Unit]
Description=Foo
[Socket]
MessageQueueMessageSize=0

