# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SyslogLevel_crit.target
[Unit]
Description=Foo
[Socket]
SyslogLevel=crit

