# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_Delegate_blkio.target
[Unit]
Description=Foo
[Socket]
Delegate=blkio

