# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_CPUSchedulingResetOnFork_no.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingResetOnFork=no

