# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KeepAlive_yes.target
[Unit]
Description=Foo
[Socket]
KeepAlive=yes

