# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENETRESET.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENETRESET

