# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_MemoryHigh_0T.target
[Unit]
Description=Foo
[Socket]
MemoryHigh=0T

