# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGSTKFLT.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGSTKFLT

