# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_BlockIOAccounting_no.target
[Unit]
Description=Foo
[Socket]
BlockIOAccounting=no

