# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_UtmpIdentifier_abc.target
[Unit]
Description=Foo
[Socket]
UtmpIdentifier=abc

