# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_IOSchedulingClass_best-effort.target
[Unit]
Description=Foo
[Socket]
IOSchedulingClass=best-effort

