# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_2.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=2

