# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_7seconds.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=7seconds

