# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TimeoutSec_26months.target
[Unit]
Description=Foo
[Socket]
TimeoutSec=26months

