# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitMSGQUEUE_infinity.target
[Unit]
Description=Foo
[Socket]
LimitMSGQUEUE=infinity

