# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitIntervalSec_11y.target
[Unit]
Description=Foo
[Socket]
TriggerLimitIntervalSec=11y

