# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGSTOP.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGSTOP

