# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SyslogLevel_info.target
[Unit]
Description=Foo
[Socket]
SyslogLevel=info

