# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_MessageQueueMaxMessages_9999999999999.target
[Unit]
Description=Foo
[Socket]
MessageQueueMaxMessages=9999999999999

