# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_MemoryAccounting_false.target
[Unit]
Description=Foo
[Socket]
MemoryAccounting=false

