# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_Transparent_on.target
[Unit]
Description=Foo
[Socket]
Transparent=on

