# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGUSR1.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGUSR1

