# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ELIBEXEC.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ELIBEXEC

