# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SyslogLevelPrefix_true.target
[Unit]
Description=Foo
[Socket]
SyslogLevelPrefix=true

