# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_26months.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=26months

