# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SendBuffer_0B.target
[Unit]
Description=Foo
[Socket]
SendBuffer=0B

