# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ECHRNG.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ECHRNG

