# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EDEADLOCK.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EDEADLOCK

