# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_15.target
[Unit]
Description=Foo
[Socket]
KillSignal=15

