# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EACCES.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EACCES

