# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_MemoryDenyWriteExecute_yes.target
[Unit]
Description=Foo
[Socket]
MemoryDenyWriteExecute=yes

