# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGUSR1.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGUSR1

