# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TTYVHangup_true.target
[Unit]
Description=Foo
[Socket]
TTYVHangup=true

