# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_PassSecurity_yes.target
[Unit]
Description=Foo
[Socket]
PassSecurity=yes

