# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_15min.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=15min

