# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGSYS.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGSYS

