# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TTYVHangup_no.target
[Unit]
Description=Foo
[Socket]
TTYVHangup=no

