# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_0sec.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=0sec

