# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_DisableControllers_yes.target
[Unit]
Description=Foo
[Socket]
DisableControllers=yes

