# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_DisableControllers_cpuacct.target
[Unit]
Description=Foo
[Socket]
DisableControllers=cpuacct

