# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_DirectoryMode_644.target
[Unit]
Description=Foo
[Socket]
DirectoryMode=644

