# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_PrivateTmp_yes.target
[Unit]
Description=Foo
[Socket]
PrivateTmp=yes

