# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_ListenSpecial_--tmp.target
[Unit]
Description=Foo
[Socket]
ListenSpecial=/tmp

