# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitMEMLOCK_9999999999999.target
[Unit]
Description=Foo
[Socket]
LimitMEMLOCK=9999999999999

