# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SendSIGKILL_false.target
[Unit]
Description=Foo
[Socket]
SendSIGKILL=false

