# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_36.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=36

