# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_RemoveOnStop_1.target
[Unit]
Description=Foo
[Socket]
RemoveOnStop=1

