# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_56.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=56

