# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_PrivateTmp_off.target
[Unit]
Description=Foo
[Socket]
PrivateTmp=off

