# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TimeoutSec_22hr.target
[Unit]
Description=Foo
[Socket]
TimeoutSec=22hr

