# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TimeoutSec_5hr.target
[Unit]
Description=Foo
[Socket]
TimeoutSec=5hr

