# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TTYVHangup_on.target
[Unit]
Description=Foo
[Socket]
TTYVHangup=on

