# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SecureBits_noroot-locked.target
[Unit]
Description=Foo
[Socket]
SecureBits=noroot-locked

