# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TTYPath_--dev--bar.target
[Unit]
Description=Foo
[Socket]
TTYPath=/dev/bar

