# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_MemoryDenyWriteExecute_true.target
[Unit]
Description=Foo
[Socket]
MemoryDenyWriteExecute=true

