# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_ReusePort_off.target
[Unit]
Description=Foo
[Socket]
ReusePort=off

