# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOTCONN.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOTCONN

