# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SyslogFacility_user.target
[Unit]
Description=Foo
[Socket]
SyslogFacility=user

