# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitIntervalSec_12second.target
[Unit]
Description=Foo
[Socket]
TriggerLimitIntervalSec=12second

