# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_24years.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=24years

