# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_23m.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=23m

