# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_Transparent_1.target
[Unit]
Description=Foo
[Socket]
Transparent=1

