# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_RootImage_--bin.target
[Unit]
Description=Foo
[Socket]
RootImage=/bin

