# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitAS_infinity.target
[Unit]
Description=Foo
[Socket]
LimitAS=infinity

