# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitSTACK_infinity.target
[Unit]
Description=Foo
[Socket]
LimitSTACK=infinity

