# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ELOOP.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ELOOP

