# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_ReceiveBuffer_0K.target
[Unit]
Description=Foo
[Socket]
ReceiveBuffer=0K

