# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_StandardError_kmsg.target
[Unit]
Description=Foo
[Socket]
StandardError=kmsg

