# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TTYPath_--dev--foo.target
[Unit]
Description=Foo
[Socket]
TTYPath=/dev/foo

