# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGPWR.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGPWR

