# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGPROF.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGPROF

