# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_RestrictSUIDSGID_on.target
[Unit]
Description=Foo
[Socket]
RestrictSUIDSGID=on

