# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_StandardOutput_journal.target
[Unit]
Description=Foo
[Socket]
StandardOutput=journal

