# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SmackProcessLabel_abc.target
[Unit]
Description=Foo
[Socket]
SmackProcessLabel=abc

