# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KeepAlive_no.target
[Unit]
Description=Foo
[Socket]
KeepAlive=no

