# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_StandardError_journal.target
[Unit]
Description=Foo
[Socket]
StandardError=journal

