# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_12second.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=12second

