# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENODEV.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENODEV

