# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_Accept_1.target
[Unit]
Description=Foo
[Socket]
Accept=1

