# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGPROF.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGPROF

