# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGPOLL.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGPOLL

