# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_ProtectSystem_full.target
[Unit]
Description=Foo
[Socket]
ProtectSystem=full

