# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EL3RST.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EL3RST

