# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGHUP.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGHUP

