# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_11m.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=11m

