# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGQUIT.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGQUIT

