# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SendSIGKILL_on.target
[Unit]
Description=Foo
[Socket]
SendSIGKILL=on

