# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SendSIGHUP_off.target
[Unit]
Description=Foo
[Socket]
SendSIGHUP=off

