# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_StartupIOWeight_10000.target
[Unit]
Description=Foo
[Socket]
StartupIOWeight=10000

