# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EHOSTUNREACH.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EHOSTUNREACH

