# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_DisableControllers_memory.target
[Unit]
Description=Foo
[Socket]
DisableControllers=memory

