# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_Accept_no.target
[Unit]
Description=Foo
[Socket]
Accept=no

