# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_BlockIOWriteBandwidthMax_--dev--bar%200G.target
[Unit]
Description=Foo
[Socket]
BlockIOWriteBandwidthMax=/dev/bar 0G

