# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TimeoutSec_28h.target
[Unit]
Description=Foo
[Socket]
TimeoutSec=28h

