# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_MemoryHigh_9999999999999M.target
[Unit]
Description=Foo
[Socket]
MemoryHigh=9999999999999M

