# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitNPROC_9999999999999.target
[Unit]
Description=Foo
[Socket]
LimitNPROC=9999999999999

