# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SendSIGKILL_yes.target
[Unit]
Description=Foo
[Socket]
SendSIGKILL=yes

