# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_2.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=2

