# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_CPUSchedulingPolicy_rr.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingPolicy=rr

