# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SyslogFacility_local4.target
[Unit]
Description=Foo
[Socket]
SyslogFacility=local4

