# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitIntervalSec_9month.target
[Unit]
Description=Foo
[Socket]
TriggerLimitIntervalSec=9month

