# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SendSIGHUP_on.target
[Unit]
Description=Foo
[Socket]
SendSIGHUP=on

