# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGPOLL.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGPOLL

