# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TimeoutSec_27y.target
[Unit]
Description=Foo
[Socket]
TimeoutSec=27y

