# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_RemoveIPC_false.target
[Unit]
Description=Foo
[Socket]
RemoveIPC=false

