# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_PrivateDevices_1.target
[Unit]
Description=Foo
[Socket]
PrivateDevices=1

