# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillMode_mixed.target
[Unit]
Description=Foo
[Socket]
KillMode=mixed

