# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitIntervalSec_13minute.target
[Unit]
Description=Foo
[Socket]
TriggerLimitIntervalSec=13minute

