# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitCORE_0.target
[Unit]
Description=Foo
[Socket]
LimitCORE=0

