# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EUNATCH.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EUNATCH

