# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_MessageQueueMessageSize_9999999999999.target
[Unit]
Description=Foo
[Socket]
MessageQueueMessageSize=9999999999999

