# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitCORE_infinity.target
[Unit]
Description=Foo
[Socket]
LimitCORE=infinity

