# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_StartupCPUWeight_10000.target
[Unit]
Description=Foo
[Socket]
StartupCPUWeight=10000

