# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_5hr.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=5hr

