# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_RestrictRealtime_off.target
[Unit]
Description=Foo
[Socket]
RestrictRealtime=off

