# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOMEDIUM.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOMEDIUM

