# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EPROTOTYPE.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EPROTOTYPE

