# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOENT.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOENT

