# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_DATAERR.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=DATAERR

