# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SyslogFacility_lpr.target
[Unit]
Description=Foo
[Socket]
SyslogFacility=lpr

