# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SyslogLevel_warning.target
[Unit]
Description=Foo
[Socket]
SyslogLevel=warning

