# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ETIMEDOUT.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ETIMEDOUT

