# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_DirectoryMode_006.target
[Unit]
Description=Foo
[Socket]
DirectoryMode=006

