# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGTSTP.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGTSTP

