# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SELinuxContextFromNet_on.target
[Unit]
Description=Foo
[Socket]
SELinuxContextFromNet=on

