# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_InaccessableDirectories_abc.target
[Unit]
Description=Foo
[Socket]
InaccessableDirectories=abc

