# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ERESTART.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ERESTART

