# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_0.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=0

