# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGTRAP.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGTRAP

