# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGVTALRM.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGVTALRM

