# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ESPIPE.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ESPIPE

