# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_TTYVTDisallocate_no.target
[Unit]
Description=Foo
[Socket]
TTYVTDisallocate=no

