# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitRTPRIO_0.target
[Unit]
Description=Foo
[Socket]
LimitRTPRIO=0

