# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_35.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=35

