# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_FileDescriptorName_abc.target
[Unit]
Description=Foo
[Socket]
FileDescriptorName=abc

