# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_BlockIODeviceWeight_--dev--bar%2010000.target
[Unit]
Description=Foo
[Socket]
BlockIODeviceWeight=/dev/bar 10000

