# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_LimitSIGPENDING_infinity.target
[Unit]
Description=Foo
[Socket]
LimitSIGPENDING=infinity

