# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_MemoryAccounting_off.target
[Unit]
Description=Foo
[Socket]
MemoryAccounting=off

