# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_Transparent_true.target
[Unit]
Description=Foo
[Socket]
Transparent=true

