# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_RestrictSUIDSGID_yes.target
[Unit]
Description=Foo
[Socket]
RestrictSUIDSGID=yes

