# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_BlockIOReadBandwidthMax_--dev--foo%200T.target
[Unit]
Description=Foo
[Socket]
BlockIOReadBandwidthMax=/dev/foo 0T

