# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_CPUSchedulingResetOnFork_yes.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingResetOnFork=yes

