# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_Mark_9999999999999.target
[Unit]
Description=Foo
[Socket]
Mark=9999999999999

