# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOTSOCK.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOTSOCK

