# systemdlint/tests/auto/InvalidValue/2.31/good/Socket_UMask_0.target
[Unit]
Description=Foo
[Socket]
UMask=0

