# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_RootImage_--bin.target
[Unit]
Description=Foo
[Swap]
RootImage=/bin

