# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGXCPU.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGXCPU

