# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_11y.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=11y

