# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_DATAERR.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=DATAERR

