# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_AppArmorProfile_abc.target
[Unit]
Description=Foo
[Swap]
AppArmorProfile=abc

