# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ESTALE.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ESTALE

