# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitNPROC_0.target
[Unit]
Description=Foo
[Swap]
LimitNPROC=0

