# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_StandardOutput_syslog%2Bconsole.target
[Unit]
Description=Foo
[Swap]
StandardOutput=syslog+console

