# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitRTTIME_infinity.target
[Unit]
Description=Foo
[Swap]
LimitRTTIME=infinity

