# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_IgnoreSIGPIPE_yes.target
[Unit]
Description=Foo
[Swap]
IgnoreSIGPIPE=yes

