# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_ProtectSystem_off.target
[Unit]
Description=Foo
[Swap]
ProtectSystem=off

