# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EMSGSIZE.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EMSGSIZE

