# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_E2BIG.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=E2BIG

