# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_ProtectHome_true.target
[Unit]
Description=Foo
[Swap]
ProtectHome=true

