# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_MountFlags_slave.target
[Unit]
Description=Foo
[Swap]
MountFlags=slave

