# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGQUIT.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGQUIT

