# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_ReadOnlyDirectories_abc.target
[Unit]
Description=Foo
[Swap]
ReadOnlyDirectories=abc

