# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitSTACK_0.target
[Unit]
Description=Foo
[Swap]
LimitSTACK=0

