# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_BlockIOWriteBandwidthMax_--dev--bar%200T.target
[Unit]
Description=Foo
[Swap]
BlockIOWriteBandwidthMax=/dev/bar 0T

