# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_AmbientCapabilities_abc.target
[Unit]
Description=Foo
[Swap]
AmbientCapabilities=abc

