# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGCHLD.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGCHLD

