# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGILL.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGILL

