# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_BlockIOAccounting_yes.target
[Unit]
Description=Foo
[Swap]
BlockIOAccounting=yes

