# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SyslogLevelPrefix_no.target
[Unit]
Description=Foo
[Swap]
SyslogLevelPrefix=no

