# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGABRT.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGABRT

