# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_DevicePolicy_strict.target
[Unit]
Description=Foo
[Swap]
DevicePolicy=strict

