# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_65.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=65

