# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EAFNOSUPPORT.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EAFNOSUPPORT

