# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EREMOTEIO.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EREMOTEIO

