# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitLOCKS_infinity.target
[Unit]
Description=Foo
[Swap]
LimitLOCKS=infinity

