# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_IOAccounting_1.target
[Unit]
Description=Foo
[Swap]
IOAccounting=1

