# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_RestrictRealtime_yes.target
[Unit]
Description=Foo
[Swap]
RestrictRealtime=yes

