# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_IOSchedulingClass_none.target
[Unit]
Description=Foo
[Swap]
IOSchedulingClass=none

