# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SELinuxContext_abc.target
[Unit]
Description=Foo
[Swap]
SELinuxContext=abc

