# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_MemoryMax_0G.target
[Unit]
Description=Foo
[Swap]
MemoryMax=0G

