# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EL3RST.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EL3RST

