# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_Delegate_yes.target
[Unit]
Description=Foo
[Swap]
Delegate=yes

