# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_MountAPIVFS_no.target
[Unit]
Description=Foo
[Swap]
MountAPIVFS=no

