# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGSEGV.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGSEGV

