# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGIO.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGIO

