# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGTTIN.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGTTIN

