# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitCORE_infinity.target
[Unit]
Description=Foo
[Swap]
LimitCORE=infinity

