# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_MemoryLimit_infinity.target
[Unit]
Description=Foo
[Swap]
MemoryLimit=infinity

