# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_ProtectSystem_on.target
[Unit]
Description=Foo
[Swap]
ProtectSystem=on

