# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_MountAPIVFS_0.target
[Unit]
Description=Foo
[Swap]
MountAPIVFS=0

