# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SyslogLevel_debug.target
[Unit]
Description=Foo
[Swap]
SyslogLevel=debug

