# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGTERM.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGTERM

