# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_Priority_9999999999999.target
[Unit]
Description=Foo
[Swap]
Priority=9999999999999

