# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitRTTIME_0.target
[Unit]
Description=Foo
[Swap]
LimitRTTIME=0

