# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_41.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=41

