# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_33.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=33

