# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EOPNOTSUPP.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EOPNOTSUPP

