# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOMSG.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOMSG

