# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGEMT.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGEMT

