# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EOVERFLOW.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EOVERFLOW

