# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EDOM.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EDOM

