# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_IOAccounting_yes.target
[Unit]
Description=Foo
[Swap]
IOAccounting=yes

