# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_6days.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=6days

