# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_DevicePolicy_auto.target
[Unit]
Description=Foo
[Swap]
DevicePolicy=auto

