# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOSPC.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOSPC

