# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitDATA_0.target
[Unit]
Description=Foo
[Swap]
LimitDATA=0

