# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGHUP.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGHUP

