# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EINVAL.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EINVAL

