# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOMEDIUM.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOMEDIUM

