# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitMSGQUEUE_infinity.target
[Unit]
Description=Foo
[Swap]
LimitMSGQUEUE=infinity

