# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_User_iamauser.target
[Unit]
Description=Foo
[Swap]
User=iamauser

