# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EADDRINUSE.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EADDRINUSE

