# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_UtmpIdentifier_abc.target
[Unit]
Description=Foo
[Swap]
UtmpIdentifier=abc

