# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_2us.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=2us

