# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_EnvironmentFile_--dev--.target
[Unit]
Description=Foo
[Swap]
EnvironmentFile=/dev/

