# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_RestrictSUIDSGID_yes.target
[Unit]
Description=Foo
[Swap]
RestrictSUIDSGID=yes

