# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_Personality_x86.target
[Unit]
Description=Foo
[Swap]
Personality=x86

