# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_StartupCPUWeight_1.target
[Unit]
Description=Foo
[Swap]
StartupCPUWeight=1

