# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGSEGV.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGSEGV

