# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitMEMLOCK_9999999999999.target
[Unit]
Description=Foo
[Swap]
LimitMEMLOCK=9999999999999

