# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_PrivateTmp_true.target
[Unit]
Description=Foo
[Swap]
PrivateTmp=true

