# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_54.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=54

