# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ECONNABORTED.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ECONNABORTED

