# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGXFSZ.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGXFSZ

