# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_Delegate_memory.target
[Unit]
Description=Foo
[Swap]
Delegate=memory

