# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SyslogFacility_cron.target
[Unit]
Description=Foo
[Swap]
SyslogFacility=cron

