# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_UMask_9999999999999.target
[Unit]
Description=Foo
[Swap]
UMask=9999999999999

