# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGTRAP.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGTRAP

