# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_16week.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=16week

