# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGPIPE.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGPIPE

