# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGFPE.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGFPE

