# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_CPUSchedulingResetOnFork_yes.target
[Unit]
Description=Foo
[Swap]
CPUSchedulingResetOnFork=yes

