# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TTYVTDisallocate_off.target
[Unit]
Description=Foo
[Swap]
TTYVTDisallocate=off

