# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SupplementaryGroups_metoo.target
[Unit]
Description=Foo
[Swap]
SupplementaryGroups=metoo

