# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOENT.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOENT

