# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TasksAccounting_off.target
[Unit]
Description=Foo
[Swap]
TasksAccounting=off

