# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_30.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=30

