# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_13.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=13

