# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_ProtectSystem_0.target
[Unit]
Description=Foo
[Swap]
ProtectSystem=0

