# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_KillSignal_6.target
[Unit]
Description=Foo
[Swap]
KillSignal=6

