# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGHUP.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGHUP

