# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EINPROGRESS.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EINPROGRESS

