# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_ReadOnlyPaths_abc.target
[Unit]
Description=Foo
[Swap]
ReadOnlyPaths=abc

