# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_MemoryAccounting_off.target
[Unit]
Description=Foo
[Swap]
MemoryAccounting=off

