# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitCORE_0.target
[Unit]
Description=Foo
[Swap]
LimitCORE=0

