# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SyslogLevel_alert.target
[Unit]
Description=Foo
[Swap]
SyslogLevel=alert

