# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_21usec.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=21usec

