# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_StandardError_socket.target
[Unit]
Description=Foo
[Swap]
StandardError=socket

