# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_15ms.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=15ms

