# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_ProtectSystem_false.target
[Unit]
Description=Foo
[Swap]
ProtectSystem=false

