# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_RootDirectory_--bin.target
[Unit]
Description=Foo
[Swap]
RootDirectory=/bin

