# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENXIO.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENXIO

