# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOSYS.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOSYS

