# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_PAMName_abc.target
[Unit]
Description=Foo
[Swap]
PAMName=abc

