# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SyslogFacility_authpriv.target
[Unit]
Description=Foo
[Swap]
SyslogFacility=authpriv

