# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_IgnoreSIGPIPE_0.target
[Unit]
Description=Foo
[Swap]
IgnoreSIGPIPE=0

