# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ESPIPE.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ESPIPE

