# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SecureBits_noroot-locked.target
[Unit]
Description=Foo
[Swap]
SecureBits=noroot-locked

