# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitNICE_0.target
[Unit]
Description=Foo
[Swap]
LimitNICE=0

