# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TTYVTDisallocate_0.target
[Unit]
Description=Foo
[Swap]
TTYVTDisallocate=0

