# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_9month.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=9month

