# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ETIMEDOUT.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ETIMEDOUT

