# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGSTOP.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGSTOP

