# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TTYVTDisallocate_yes.target
[Unit]
Description=Foo
[Swap]
TTYVTDisallocate=yes

