# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_17second.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=17second

