# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TTYVTDisallocate_false.target
[Unit]
Description=Foo
[Swap]
TTYVTDisallocate=false

