# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_Personality_ppc-le.target
[Unit]
Description=Foo
[Swap]
Personality=ppc-le

