# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EMFILE.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EMFILE

