# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ECHILD.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ECHILD

