# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGQUIT.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGQUIT

