# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EADDRNOTAVAIL.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EADDRNOTAVAIL

