# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOTUNIQ.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOTUNIQ

