# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_StandardOutput_null.target
[Unit]
Description=Foo
[Swap]
StandardOutput=null

