# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SyslogFacility_kern.target
[Unit]
Description=Foo
[Swap]
SyslogFacility=kern

