# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_ProtectSystem_full.target
[Unit]
Description=Foo
[Swap]
ProtectSystem=full

