# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_ProtectSystem_yes.target
[Unit]
Description=Foo
[Swap]
ProtectSystem=yes

