# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallArchitectures_abc.target
[Unit]
Description=Foo
[Swap]
SystemCallArchitectures=abc

