# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WorkingDirectory_--bin.target
[Unit]
Description=Foo
[Swap]
WorkingDirectory=/bin

