# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGWINCH.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGWINCH

