# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_26months.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=26months

