# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EDESTADDRREQ.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EDESTADDRREQ

