# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_ProtectHome_tmpfs.target
[Unit]
Description=Foo
[Swap]
ProtectHome=tmpfs

