# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SyslogFacility_daemon.target
[Unit]
Description=Foo
[Swap]
SyslogFacility=daemon

