# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SyslogLevel_info.target
[Unit]
Description=Foo
[Swap]
SyslogLevel=info

