# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGSEGV.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGSEGV

