# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EIDRM.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EIDRM

