# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_RestrictRealtime_false.target
[Unit]
Description=Foo
[Swap]
RestrictRealtime=false

