# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_RemoveIPC_off.target
[Unit]
Description=Foo
[Swap]
RemoveIPC=off

