# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_MemoryLimit_0B.target
[Unit]
Description=Foo
[Swap]
MemoryLimit=0B

