# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitCORE_9999999999999.target
[Unit]
Description=Foo
[Swap]
LimitCORE=9999999999999

