# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WorkingDirectory_--usr--sbin.target
[Unit]
Description=Foo
[Swap]
WorkingDirectory=/usr/sbin

