# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_CPUQuota_0%25.target
[Unit]
Description=Foo
[Swap]
CPUQuota=0%

