# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_55.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=55

