# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENETDOWN.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENETDOWN

