# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SyslogLevel_err.target
[Unit]
Description=Foo
[Swap]
SyslogLevel=err

