# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_DisableControllers_cpu.target
[Unit]
Description=Foo
[Swap]
DisableControllers=cpu

