# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_PrivateNetwork_false.target
[Unit]
Description=Foo
[Swap]
PrivateNetwork=false

