# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitAS_9999999999999.target
[Unit]
Description=Foo
[Swap]
LimitAS=9999999999999

