# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ELOOP.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ELOOP

