# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGABRT.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGABRT

