# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_20w.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=20w

