# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_3hours.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=3hours

