# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGALRM.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGALRM

