# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_Delegate_cpu.target
[Unit]
Description=Foo
[Swap]
Delegate=cpu

