# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_BlockStartupIOWeight_1.target
[Unit]
Description=Foo
[Swap]
BlockStartupIOWeight=1

