# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGUSR1.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGUSR1

