# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitNPROC_9999999999999.target
[Unit]
Description=Foo
[Swap]
LimitNPROC=9999999999999

