# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_ProtectSystem_strict.target
[Unit]
Description=Foo
[Swap]
ProtectSystem=strict

