# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_ProtectSystem_true.target
[Unit]
Description=Foo
[Swap]
ProtectSystem=true

