# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_IOSchedulingPriority_7.target
[Unit]
Description=Foo
[Swap]
IOSchedulingPriority=7

