# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_IgnoreSIGPIPE_true.target
[Unit]
Description=Foo
[Swap]
IgnoreSIGPIPE=true

