# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EKEYREJECTED.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EKEYREJECTED

