# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_PrivateTmp_off.target
[Unit]
Description=Foo
[Swap]
PrivateTmp=off

