# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SyslogFacility_news.target
[Unit]
Description=Foo
[Swap]
SyslogFacility=news

