# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EREMCHG.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EREMCHG

