# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_MemoryMax_0M.target
[Unit]
Description=Foo
[Swap]
MemoryMax=0M

