# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENAMETOOLONG.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENAMETOOLONG

