# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_IgnoreSIGPIPE_false.target
[Unit]
Description=Foo
[Swap]
IgnoreSIGPIPE=false

