# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_IgnoreSIGPIPE_1.target
[Unit]
Description=Foo
[Swap]
IgnoreSIGPIPE=1

