# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGCONT.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGCONT

