# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_PrivateTmp_no.target
[Unit]
Description=Foo
[Swap]
PrivateTmp=no

