# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGPOLL.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGPOLL

