# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_24.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=24

