# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitRTPRIO_infinity.target
[Unit]
Description=Foo
[Swap]
LimitRTPRIO=infinity

