# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGSTOP.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGSTOP

