# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGVTALRM.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGVTALRM

