# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EACCES.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EACCES

