# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENONET.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENONET

