# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SyslogLevelPrefix_on.target
[Unit]
Description=Foo
[Swap]
SyslogLevelPrefix=on

