# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGILL.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGILL

