# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_MountAPIVFS_yes.target
[Unit]
Description=Foo
[Swap]
MountAPIVFS=yes

