# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TTYReset_off.target
[Unit]
Description=Foo
[Swap]
TTYReset=off

