# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGCLD.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGCLD

