# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGCLD.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGCLD

