# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EISCONN.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EISCONN

