# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGIO.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGIO

