# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOLINK.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOLINK

