# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EWOULDBLOCK.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EWOULDBLOCK

