# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EINTR.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EINTR

