# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_5hour.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=5hour

