# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EHWPOISON.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EHWPOISON

