# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGEMT.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGEMT

