# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGIO.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGIO

