# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TasksAccounting_1.target
[Unit]
Description=Foo
[Swap]
TasksAccounting=1

