# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_MemoryAccounting_yes.target
[Unit]
Description=Foo
[Swap]
MemoryAccounting=yes

