# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SendSIGHUP_false.target
[Unit]
Description=Foo
[Swap]
SendSIGHUP=false

