# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_CPUSchedulingPolicy_batch.target
[Unit]
Description=Foo
[Swap]
CPUSchedulingPolicy=batch

