# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_BlockStartupIOWeight_10000.target
[Unit]
Description=Foo
[Swap]
BlockStartupIOWeight=10000

