# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_BlockIODeviceWeight_--dev--bar%2010000.target
[Unit]
Description=Foo
[Swap]
BlockIODeviceWeight=/dev/bar 10000

