# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_PrivateTmp_1.target
[Unit]
Description=Foo
[Swap]
PrivateTmp=1

