# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGTRAP.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGTRAP

