# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_KillSignal_8.target
[Unit]
Description=Foo
[Swap]
KillSignal=8

