# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_Delegate_cpuacct.target
[Unit]
Description=Foo
[Swap]
Delegate=cpuacct

