# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitSTACK_infinity.target
[Unit]
Description=Foo
[Swap]
LimitSTACK=infinity

