# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGPROF.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGPROF

