# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_ReadWritePaths_abc.target
[Unit]
Description=Foo
[Swap]
ReadWritePaths=abc

