# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_StartupIOWeight_10000.target
[Unit]
Description=Foo
[Swap]
StartupIOWeight=10000

