# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOTCONN.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOTCONN

