# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EKEYEXPIRED.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EKEYEXPIRED

