# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGTSTP.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGTSTP

