# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ESTRPIPE.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ESTRPIPE

