# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_Priority_0.target
[Unit]
Description=Foo
[Swap]
Priority=0

