# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SyslogIdentifier_abc.target
[Unit]
Description=Foo
[Swap]
SyslogIdentifier=abc

