# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_BlockIOAccounting_0.target
[Unit]
Description=Foo
[Swap]
BlockIOAccounting=0

