# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_1.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=1

