# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGILL.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGILL

