# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_25s.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=25s

