# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_2.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=2

