# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EPROTOTYPE.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EPROTOTYPE

