# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SecureBits_no-setuid-fixup.target
[Unit]
Description=Foo
[Swap]
SecureBits=no-setuid-fixup

