# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ERANGE.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ERANGE

