# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_IgnoreSIGPIPE_no.target
[Unit]
Description=Foo
[Swap]
IgnoreSIGPIPE=no

