# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_StandardOutput_inherit.target
[Unit]
Description=Foo
[Swap]
StandardOutput=inherit

