# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOTDIR.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOTDIR

