# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENODEV.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENODEV

