# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitDATA_9999999999999.target
[Unit]
Description=Foo
[Swap]
LimitDATA=9999999999999

