# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGSYS.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGSYS

