# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TTYVHangup_off.target
[Unit]
Description=Foo
[Swap]
TTYVHangup=off

