# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_StartupCPUShares_262144.target
[Unit]
Description=Foo
[Swap]
StartupCPUShares=262144

