# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGCHLD.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGCHLD

