# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EISDIR.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EISDIR

