# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGKILL.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGKILL

