# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitCPU_0.target
[Unit]
Description=Foo
[Swap]
LimitCPU=0

