# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TasksAccounting_no.target
[Unit]
Description=Foo
[Swap]
TasksAccounting=no

