# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_11m.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=11m

