# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGALRM.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGALRM

