# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_23.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=23

