# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGFPE.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGFPE

