# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGPROF.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGPROF

