# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_0.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=0

