# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGINFO.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGINFO

