# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_PrivateDevices_no.target
[Unit]
Description=Foo
[Swap]
PrivateDevices=no

