# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitFSIZE_9999999999999.target
[Unit]
Description=Foo
[Swap]
LimitFSIZE=9999999999999

