# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SyslogFacility_lpr.target
[Unit]
Description=Foo
[Swap]
SyslogFacility=lpr

