# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_MountFlags_private.target
[Unit]
Description=Foo
[Swap]
MountFlags=private

