# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOTBLK.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOTBLK

