# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOTTY.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOTTY

