# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGSTKFLT.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGSTKFLT

