# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TimerSlackNSec_9999999999999.target
[Unit]
Description=Foo
[Swap]
TimerSlackNSec=9999999999999

