# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_Delegate_blkio.target
[Unit]
Description=Foo
[Swap]
Delegate=blkio

