# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGHUP.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGHUP

