# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_IOWriteBandwidthMax_abc.target
[Unit]
Description=Foo
[Swap]
IOWriteBandwidthMax=abc

