# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitNICE_infinity.target
[Unit]
Description=Foo
[Swap]
LimitNICE=infinity

