# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_MemorySwapMax_9999999999999G.target
[Unit]
Description=Foo
[Swap]
MemorySwapMax=9999999999999G

