# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_What_abc.target
[Unit]
Description=Foo
[Swap]
What=abc

