# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_19d.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=19d

