# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_MemoryHigh_0B.target
[Unit]
Description=Foo
[Swap]
MemoryHigh=0B

