# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_DATAERR.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=DATAERR

