# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EDQUOT.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EDQUOT

