# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_PrivateDevices_true.target
[Unit]
Description=Foo
[Swap]
PrivateDevices=true

