# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_BlockIOWriteBandwidthMax_--dev--foo%200B.target
[Unit]
Description=Foo
[Swap]
BlockIOWriteBandwidthMax=/dev/foo 0B

