# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SmackProcessLabel_abc.target
[Unit]
Description=Foo
[Swap]
SmackProcessLabel=abc

