# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGURG.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGURG

