# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_IOAccounting_off.target
[Unit]
Description=Foo
[Swap]
IOAccounting=off

