# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ELIBSCN.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ELIBSCN

