# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_MemoryLow_9999999999999K.target
[Unit]
Description=Foo
[Swap]
MemoryLow=9999999999999K

