# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_CPUSchedulingResetOnFork_off.target
[Unit]
Description=Foo
[Swap]
CPUSchedulingResetOnFork=off

