# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_IOSchedulingClass_idle.target
[Unit]
Description=Foo
[Swap]
IOSchedulingClass=idle

