# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGCLD.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGCLD

