# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EXDEV.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EXDEV

