# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_MemoryDenyWriteExecute_yes.target
[Unit]
Description=Foo
[Swap]
MemoryDenyWriteExecute=yes

