# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_MemoryDenyWriteExecute_0.target
[Unit]
Description=Foo
[Swap]
MemoryDenyWriteExecute=0

