# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_StandardOutput_journal%2Bconsole.target
[Unit]
Description=Foo
[Swap]
StandardOutput=journal+console

