# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TasksMax_1.target
[Unit]
Description=Foo
[Swap]
TasksMax=1

