# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitNOFILE_9999999999999.target
[Unit]
Description=Foo
[Swap]
LimitNOFILE=9999999999999

