# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_6.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=6

