# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EDEADLK.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EDEADLK

