# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOTSOCK.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOTSOCK

