# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOMEM.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOMEM

