# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TTYVHangup_no.target
[Unit]
Description=Foo
[Swap]
TTYVHangup=no

