# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitMSGQUEUE_0.target
[Unit]
Description=Foo
[Swap]
LimitMSGQUEUE=0

