# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGPOLL.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGPOLL

