# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_CPUSchedulingPriority_1.target
[Unit]
Description=Foo
[Swap]
CPUSchedulingPriority=1

