# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_RootImage_--usr--sbin.target
[Unit]
Description=Foo
[Swap]
RootImage=/usr/sbin

