# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EREMOTE.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EREMOTE

