# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_StandardError_syslog.target
[Unit]
Description=Foo
[Swap]
StandardError=syslog

