# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_IOReadBandwidthMax_abc.target
[Unit]
Description=Foo
[Swap]
IOReadBandwidthMax=abc

