# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_1%C2%B5s.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=1µs

