# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TTYVTDisallocate_no.target
[Unit]
Description=Foo
[Swap]
TTYVTDisallocate=no

