# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TasksMax_infinity.target
[Unit]
Description=Foo
[Swap]
TasksMax=infinity

