# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EPFNOSUPPORT.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EPFNOSUPPORT

