# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitMEMLOCK_0.target
[Unit]
Description=Foo
[Swap]
LimitMEMLOCK=0

