# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EUSERS.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EUSERS

