# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_StandardError_inherit.target
[Unit]
Description=Foo
[Swap]
StandardError=inherit

