# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitNOFILE_infinity.target
[Unit]
Description=Foo
[Swap]
LimitNOFILE=infinity

