# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_RemoveIPC_false.target
[Unit]
Description=Foo
[Swap]
RemoveIPC=false

