# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_27msec.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=27msec

