# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_StandardError_kmsg%2Bconsole.target
[Unit]
Description=Foo
[Swap]
StandardError=kmsg+console

