# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_BlockIOReadBandwidthMax_--dev--bar%209999999999999M.target
[Unit]
Description=Foo
[Swap]
BlockIOReadBandwidthMax=/dev/bar 9999999999999M

