# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_UtmpMode_init.target
[Unit]
Description=Foo
[Swap]
UtmpMode=init

