# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitSIGPENDING_0.target
[Unit]
Description=Foo
[Swap]
LimitSIGPENDING=0

