# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EL2HLT.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EL2HLT

