# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SendSIGKILL_no.target
[Unit]
Description=Foo
[Swap]
SendSIGKILL=no

