# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_LimitRSS_9999999999999.target
[Unit]
Description=Foo
[Swap]
LimitRSS=9999999999999

