# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EBADRQC.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EBADRQC

