# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_RestrictRealtime_true.target
[Unit]
Description=Foo
[Swap]
RestrictRealtime=true

