# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ECONNRESET.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ECONNRESET

