# systemdlint/tests/auto/InvalidValue/2.31/good/Swap_PrivateDevices_0.target
[Unit]
Description=Foo
[Swap]
PrivateDevices=0

