# systemdlint/tests/auto/InvalidValue/2.31/good/Tap_PacketInfo_yes.target
[Unit]
Description=Foo
[Tap]
PacketInfo=yes

