# systemdlint/tests/auto/InvalidValue/2.31/good/Tap_PacketInfo_0.target
[Unit]
Description=Foo
[Tap]
PacketInfo=0

