# systemdlint/tests/auto/InvalidValue/2.31/good/Tap_OneQueue_1.target
[Unit]
Description=Foo
[Tap]
OneQueue=1

