# systemdlint/tests/auto/InvalidValue/2.31/good/Tap_OneQueue_0.target
[Unit]
Description=Foo
[Tap]
OneQueue=0

