# systemdlint/tests/auto/InvalidValue/2.31/good/Tap_MultiQueue_false.target
[Unit]
Description=Foo
[Tap]
MultiQueue=false

