# systemdlint/tests/auto/InvalidValue/2.31/good/Tap_VNetHeader_0.target
[Unit]
Description=Foo
[Tap]
VNetHeader=0

