# systemdlint/tests/auto/InvalidValue/2.31/good/Tap_PacketInfo_on.target
[Unit]
Description=Foo
[Tap]
PacketInfo=on

