# systemdlint/tests/auto/InvalidValue/2.31/good/Tap_MultiQueue_true.target
[Unit]
Description=Foo
[Tap]
MultiQueue=true

