# systemdlint/tests/auto/InvalidValue/2.31/good/Tap_MultiQueue_off.target
[Unit]
Description=Foo
[Tap]
MultiQueue=off

