# systemdlint/tests/auto/InvalidValue/2.31/good/Tap_OneQueue_true.target
[Unit]
Description=Foo
[Tap]
OneQueue=true

