# systemdlint/tests/auto/InvalidValue/2.31/good/Tap_MultiQueue_no.target
[Unit]
Description=Foo
[Tap]
MultiQueue=no

