# systemdlint/tests/auto/InvalidValue/2.31/good/Tap_PacketInfo_false.target
[Unit]
Description=Foo
[Tap]
PacketInfo=false

