# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_10M.target
[Unit]
Description=Foo
[Timer]
AccuracySec=10M

