# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnClockChange_1.target
[Unit]
Description=Foo
[Timer]
OnClockChange=1

