# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_26months.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=26months

