# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnStartupSec_infinity.target
[Unit]
Description=Foo
[Timer]
OnStartupSec=infinity

