# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_8minute.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=8minute

