# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_infinity.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=infinity

