# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnTimezoneChange_0.target
[Unit]
Description=Foo
[Timer]
OnTimezoneChange=0

