# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_28h.target
[Unit]
Description=Foo
[Timer]
AccuracySec=28h

