# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_10M.target
[Unit]
Description=Foo
[Timer]
OnBootSec=10M

