# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitActiveSec_27msec.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=27msec

