# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_28h.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=28h

