# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_7minutes.target
[Unit]
Description=Foo
[Timer]
AccuracySec=7minutes

