# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnClockChange_false.target
[Unit]
Description=Foo
[Timer]
OnClockChange=false

