# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_24years.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=24years

