# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_21usec.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=21usec

