# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_24w.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=24w

