# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_0sec.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=0sec

