# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_26months.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=26months

