# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_25years.target
[Unit]
Description=Foo
[Timer]
OnBootSec=25years

