# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnTimezoneChange_off.target
[Unit]
Description=Foo
[Timer]
OnTimezoneChange=off

