# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitActiveSec_0sec.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=0sec

