# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_26months.target
[Unit]
Description=Foo
[Timer]
OnBootSec=26months

