# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_4ms.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=4ms

