# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_24years.target
[Unit]
Description=Foo
[Timer]
OnBootSec=24years

