# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_17minutes.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=17minutes

