# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_infinity.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=infinity

