# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitActiveSec_5hour.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=5hour

