# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_9month.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=9month

