# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_7seconds.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=7seconds

