# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_24w.target
[Unit]
Description=Foo
[Timer]
AccuracySec=24w

