# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_16week.target
[Unit]
Description=Foo
[Timer]
AccuracySec=16week

