# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_24w.target
[Unit]
Description=Foo
[Timer]
OnBootSec=24w

