# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_20w.target
[Unit]
Description=Foo
[Timer]
OnBootSec=20w

