# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_17second.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=17second

