# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_18year.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=18year

