# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_12second.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=12second

