# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_19d.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=19d

