# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_20usec.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=20usec

