# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_14M.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=14M

