# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_10day.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=10day

