# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_12second.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=12second

