# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_21us.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=21us

