# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnStartupSec_5hr.target
[Unit]
Description=Foo
[Timer]
OnStartupSec=5hr

