# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitActiveSec_4hour.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=4hour

