# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_17second.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=17second

