# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_7seconds.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=7seconds

