# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_4ms.target
[Unit]
Description=Foo
[Timer]
AccuracySec=4ms

