# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_19d.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=19d

