# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_0sec.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=0sec

