# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_26year.target
[Unit]
Description=Foo
[Timer]
AccuracySec=26year

