# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_8weeks.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=8weeks

