# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_26months.target
[Unit]
Description=Foo
[Timer]
AccuracySec=26months

