# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_21usec.target
[Unit]
Description=Foo
[Timer]
OnBootSec=21usec

