# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_28%C2%B5s.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=28µs

