# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_WakeSystem_on.target
[Unit]
Description=Foo
[Timer]
WakeSystem=on

