# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_Persistent_true.target
[Unit]
Description=Foo
[Timer]
Persistent=true

