# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_9min.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=9min

