# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_1day.target
[Unit]
Description=Foo
[Timer]
OnBootSec=1day

