# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_21usec.target
[Unit]
Description=Foo
[Timer]
AccuracySec=21usec

