# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_26year.target
[Unit]
Description=Foo
[Timer]
OnBootSec=26year

