# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_0sec.target
[Unit]
Description=Foo
[Timer]
AccuracySec=0sec

