# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_16week.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=16week

