# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_16week.target
[Unit]
Description=Foo
[Timer]
OnBootSec=16week

