# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_12second.target
[Unit]
Description=Foo
[Timer]
OnBootSec=12second

