# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_Persistent_on.target
[Unit]
Description=Foo
[Timer]
Persistent=on

