# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_23week.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=23week

