# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_0days.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=0days

