# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_3hours.target
[Unit]
Description=Foo
[Timer]
AccuracySec=3hours

