# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_10day.target
[Unit]
Description=Foo
[Timer]
OnBootSec=10day

