# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_17minutes.target
[Unit]
Description=Foo
[Timer]
OnBootSec=17minutes

