# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_4hour.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=4hour

