# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_22hr.target
[Unit]
Description=Foo
[Timer]
OnBootSec=22hr

