# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_22hr.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=22hr

