# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_2us.target
[Unit]
Description=Foo
[Timer]
OnBootSec=2us

