# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_Persistent_off.target
[Unit]
Description=Foo
[Timer]
Persistent=off

