# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_8weeks.target
[Unit]
Description=Foo
[Timer]
OnBootSec=8weeks

