# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitActiveSec_9month.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=9month

