# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_3hours.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=3hours

