# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_WakeSystem_no.target
[Unit]
Description=Foo
[Timer]
WakeSystem=no

