# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_0sec.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=0sec

