# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitActiveSec_13months.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=13months

