# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_13minute.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=13minute

