# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_13minute.target
[Unit]
Description=Foo
[Timer]
OnBootSec=13minute

