# systemdlint/tests/auto/InvalidValue/2.31/good/Time_PollIntervalMaxSec_17.target
[Unit]
Description=Foo
[Time]
PollIntervalMaxSec=17

