# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_27msec.target
[Unit]
Description=Foo
[Timer]
AccuracySec=27msec

