# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_28%C2%B5s.target
[Unit]
Description=Foo
[Timer]
OnBootSec=28µs

