# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_6h.target
[Unit]
Description=Foo
[Timer]
OnBootSec=6h

