# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_15min.target
[Unit]
Description=Foo
[Timer]
OnBootSec=15min

