# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_9month.target
[Unit]
Description=Foo
[Timer]
AccuracySec=9month

