# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_9month.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=9month

