# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_25s.target
[Unit]
Description=Foo
[Timer]
AccuracySec=25s

