# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RemainAfterElapse_false.target
[Unit]
Description=Foo
[Timer]
RemainAfterElapse=false

