# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_20w.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=20w

