# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_infinity.target
[Unit]
Description=Foo
[Timer]
AccuracySec=infinity

