# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_22hr.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=22hr

