# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_20w.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=20w

