# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_18sec.target
[Unit]
Description=Foo
[Timer]
AccuracySec=18sec

