# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_23m.target
[Unit]
Description=Foo
[Timer]
AccuracySec=23m

