# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_17minutes.target
[Unit]
Description=Foo
[Timer]
AccuracySec=17minutes

