# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_5hr.target
[Unit]
Description=Foo
[Timer]
OnBootSec=5hr

