# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_7seconds.target
[Unit]
Description=Foo
[Timer]
AccuracySec=7seconds

