# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_1day.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=1day

