# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_6days.target
[Unit]
Description=Foo
[Timer]
AccuracySec=6days

