# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_6days.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=6days

