# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_14M.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=14M

