# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_22hr.target
[Unit]
Description=Foo
[Timer]
AccuracySec=22hr

