# systemdlint/tests/auto/InvalidValue/2.31/good/Time_PollIntervalMaxSec_9999999999999.target
[Unit]
Description=Foo
[Time]
PollIntervalMaxSec=9999999999999

