# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RemainAfterElapse_1.target
[Unit]
Description=Foo
[Timer]
RemainAfterElapse=1

