# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_11m.target
[Unit]
Description=Foo
[Timer]
OnBootSec=11m

