# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_4hour.target
[Unit]
Description=Foo
[Timer]
OnBootSec=4hour

