# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_2us.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=2us

