# systemdlint/tests/auto/InvalidValue/2.31/good/Time_PollIntervalMinSec_9999999999999.target
[Unit]
Description=Foo
[Time]
PollIntervalMinSec=9999999999999

