# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitActiveSec_23m.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=23m

