# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_13minute.target
[Unit]
Description=Foo
[Timer]
AccuracySec=13minute

