# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_12second.target
[Unit]
Description=Foo
[Timer]
AccuracySec=12second

