# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_6days.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=6days

