# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_11y.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=11y

