# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_18sec.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=18sec

