# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitActiveSec_infinity.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=infinity

