# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitActiveSec_24w.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=24w

