# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_2d.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=2d

