# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_18year.target
[Unit]
Description=Foo
[Timer]
OnBootSec=18year

