# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_15ms.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=15ms

