# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitActiveSec_16week.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=16week

