# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnClockChange_true.target
[Unit]
Description=Foo
[Timer]
OnClockChange=true

