# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_5hr.target
[Unit]
Description=Foo
[Timer]
AccuracySec=5hr

