# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_25s.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=25s

