# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_10day.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=10day

