# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RemainAfterElapse_on.target
[Unit]
Description=Foo
[Timer]
RemainAfterElapse=on

