# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_10day.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=10day

