# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_5hour.target
[Unit]
Description=Foo
[Timer]
AccuracySec=5hour

