# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_22weeks.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=22weeks

