# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_21usec.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=21usec

