# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_5hour.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=5hour

