# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_27y.target
[Unit]
Description=Foo
[Timer]
OnBootSec=27y

