# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_4ms.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=4ms

