# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_8minute.target
[Unit]
Description=Foo
[Timer]
OnBootSec=8minute

