# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_7seconds.target
[Unit]
Description=Foo
[Timer]
OnBootSec=7seconds

