# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_Persistent_yes.target
[Unit]
Description=Foo
[Timer]
Persistent=yes

