# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_WakeSystem_false.target
[Unit]
Description=Foo
[Timer]
WakeSystem=false

