# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_13minute.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=13minute

