# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_11y.target
[Unit]
Description=Foo
[Timer]
AccuracySec=11y

