# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_6days.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=6days

