# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_4ms.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=4ms

