# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitActiveSec_4ms.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=4ms

