# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_28%C2%B5s.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=28µs

