# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_24w.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=24w

