# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_18year.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=18year

