# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_7seconds.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=7seconds

