# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_19d.target
[Unit]
Description=Foo
[Timer]
AccuracySec=19d

