# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_9month.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=9month

