# systemdlint/tests/auto/InvalidValue/2.31/good/Time_PollIntervalMinSec_16.target
[Unit]
Description=Foo
[Time]
PollIntervalMinSec=16

