# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnUnitActiveSec_15min.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=15min

