# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_23m.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=23m

