# systemdlint/tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_19d.target
[Unit]
Description=Foo
[Timer]
OnBootSec=19d

