# systemdlint/tests/auto/InvalidValue/2.31/good/Tun_User_metoo.target
[Unit]
Description=Foo
[Tun]
User=metoo

