# systemdlint/tests/auto/InvalidValue/2.31/good/Tun_Group_iamagroup.target
[Unit]
Description=Foo
[Tun]
Group=iamagroup

