# systemdlint/tests/auto/InvalidValue/2.31/good/Tun_MultiQueue_0.target
[Unit]
Description=Foo
[Tun]
MultiQueue=0

