# systemdlint/tests/auto/InvalidValue/2.31/good/Tun_VNetHeader_true.target
[Unit]
Description=Foo
[Tun]
VNetHeader=true

