# systemdlint/tests/auto/InvalidValue/2.31/good/Tun_VNetHeader_1.target
[Unit]
Description=Foo
[Tun]
VNetHeader=1

