# systemdlint/tests/auto/InvalidValue/2.31/good/Tun_Group_metoo.target
[Unit]
Description=Foo
[Tun]
Group=metoo

