# systemdlint/tests/auto/InvalidValue/2.31/good/Tun_MultiQueue_off.target
[Unit]
Description=Foo
[Tun]
MultiQueue=off

