# systemdlint/tests/auto/InvalidValue/2.31/good/Tun_OneQueue_on.target
[Unit]
Description=Foo
[Tun]
OneQueue=on

