# systemdlint/tests/auto/InvalidValue/2.31/good/Tun_VNetHeader_on.target
[Unit]
Description=Foo
[Tun]
VNetHeader=on

