# systemdlint/tests/auto/InvalidValue/2.31/good/Tun_OneQueue_1.target
[Unit]
Description=Foo
[Tun]
OneQueue=1

