# systemdlint/tests/auto/InvalidValue/2.31/good/Tun_PacketInfo_0.target
[Unit]
Description=Foo
[Tun]
PacketInfo=0

