# systemdlint/tests/auto/InvalidValue/2.31/good/Tun_OneQueue_true.target
[Unit]
Description=Foo
[Tun]
OneQueue=true

