# systemdlint/tests/auto/InvalidValue/2.31/good/Tun_OneQueue_off.target
[Unit]
Description=Foo
[Tun]
OneQueue=off

