# systemdlint/tests/auto/InvalidValue/2.31/good/Tun_MultiQueue_true.target
[Unit]
Description=Foo
[Tun]
MultiQueue=true

