# systemdlint/tests/auto/InvalidValue/2.31/good/Tun_PacketInfo_true.target
[Unit]
Description=Foo
[Tun]
PacketInfo=true

