# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_OutputKey_abc.target
[Unit]
Description=Foo
[Tunnel]
OutputKey=abc

