# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_EncapsulationLimit_0.target
[Unit]
Description=Foo
[Tunnel]
EncapsulationLimit=0

