# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_ISATAP_off.target
[Unit]
Description=Foo
[Tunnel]
ISATAP=off

