# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_Remote_1.1.1.1.target
[Unit]
Description=Foo
[Tunnel]
Remote=1.1.1.1

