# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_FOUSourcePort_65534.target
[Unit]
Description=Foo
[Tunnel]
FOUSourcePort=65534

