# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_TTL_255.target
[Unit]
Description=Foo
[Tunnel]
TTL=255

