# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_EncapsulationLimit_none.target
[Unit]
Description=Foo
[Tunnel]
EncapsulationLimit=none

