# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_AllowLocalRemote_1.target
[Unit]
Description=Foo
[Tunnel]
AllowLocalRemote=1

