# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_Mode_ip6ip6.target
[Unit]
Description=Foo
[Tunnel]
Mode=ip6ip6

