# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_CopyDSCP_1.target
[Unit]
Description=Foo
[Tunnel]
CopyDSCP=1

