# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_SerializeTunneledPackets_false.target
[Unit]
Description=Foo
[Tunnel]
SerializeTunneledPackets=false

