# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_IPv6FlowLabel_abc.target
[Unit]
Description=Foo
[Tunnel]
IPv6FlowLabel=abc

