# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_CopyDSCP_true.target
[Unit]
Description=Foo
[Tunnel]
CopyDSCP=true

