# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_Independent_on.target
[Unit]
Description=Foo
[Tunnel]
Independent=on

