# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_Independent_yes.target
[Unit]
Description=Foo
[Tunnel]
Independent=yes

