# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_FooOverUDP_off.target
[Unit]
Description=Foo
[Tunnel]
FooOverUDP=off

