# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_CopyDSCP_0.target
[Unit]
Description=Foo
[Tunnel]
CopyDSCP=0

