# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_Local_255.255.255.255.target
[Unit]
Description=Foo
[Tunnel]
Local=255.255.255.255

