# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_CopyDSCP_off.target
[Unit]
Description=Foo
[Tunnel]
CopyDSCP=off

