# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_AllowLocalRemote_no.target
[Unit]
Description=Foo
[Tunnel]
AllowLocalRemote=no

