# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_ISATAP_false.target
[Unit]
Description=Foo
[Tunnel]
ISATAP=false

