# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_ISATAP_yes.target
[Unit]
Description=Foo
[Tunnel]
ISATAP=yes

