# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_DiscoverPathMTU_1.target
[Unit]
Description=Foo
[Tunnel]
DiscoverPathMTU=1

