# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_CopyDSCP_false.target
[Unit]
Description=Foo
[Tunnel]
CopyDSCP=false

