# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_FOUDestinationPort_65534.target
[Unit]
Description=Foo
[Tunnel]
FOUDestinationPort=65534

