# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_Remote_any.target
[Unit]
Description=Foo
[Tunnel]
Remote=any

