# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_ISATAP_on.target
[Unit]
Description=Foo
[Tunnel]
ISATAP=on

