# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_Local_any.target
[Unit]
Description=Foo
[Tunnel]
Local=any

