# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_Independent_1.target
[Unit]
Description=Foo
[Tunnel]
Independent=1

