# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_InputKey_abc.target
[Unit]
Description=Foo
[Tunnel]
InputKey=abc

