# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_AllowLocalRemote_on.target
[Unit]
Description=Foo
[Tunnel]
AllowLocalRemote=on

