# systemdlint/tests/auto/InvalidValue/2.31/good/Tunnel_SerializeTunneledPackets_1.target
[Unit]
Description=Foo
[Tunnel]
SerializeTunneledPackets=1

