# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_ConditionFileIsExecutable_--dev--.target
[Unit]
Description=Foo
Description=abc
ConditionFileIsExecutable=/dev/

