# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_ConditionFirstBoot_off.target
[Unit]
Description=Foo
Description=abc
ConditionFirstBoot=off

