# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_IgnoreOnIsolate_true.target
[Unit]
Description=Foo
Description=abc
IgnoreOnIsolate=true

