# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_ConditionArchitecture_abc.target
[Unit]
Description=Foo
Description=abc
ConditionArchitecture=abc

