# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_AssertKernelCommandLine_abc.target
[Unit]
Description=Foo
Description=abc
AssertKernelCommandLine=abc

