# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobRunningTimeoutSec_26months.target
[Unit]
Description=Foo
Description=abc
JobRunningTimeoutSec=26months

