# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_AssertHost_abc.target
[Unit]
Description=Foo
Description=abc
AssertHost=abc

