# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_AssertFirstBoot_1.target
[Unit]
Description=Foo
Description=abc
AssertFirstBoot=1

