# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_AssertNeedsUpdate_abc.target
[Unit]
Description=Foo
Description=abc
AssertNeedsUpdate=abc

