# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_AllowIsolate_no.target
[Unit]
Description=Foo
Description=abc
AllowIsolate=no

