# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutAction_none.target
[Unit]
Description=Foo
Description=abc
JobTimeoutAction=none

