# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobRunningTimeoutSec_17minutes.target
[Unit]
Description=Foo
Description=abc
JobRunningTimeoutSec=17minutes

