# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_OnFailureJobMode_flush.target
[Unit]
Description=Foo
Description=abc
OnFailureJobMode=flush

