# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_AssertACPower_0.target
[Unit]
Description=Foo
Description=abc
AssertACPower=0

