# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobRunningTimeoutSec_18year.target
[Unit]
Description=Foo
Description=abc
JobRunningTimeoutSec=18year

