# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_ConditionDirectoryNotEmpty_--dev--.target
[Unit]
Description=Foo
Description=abc
ConditionDirectoryNotEmpty=/dev/

