# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_StopWhenUnneeded_1.target
[Unit]
Description=Foo
Description=abc
StopWhenUnneeded=1

