# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobRunningTimeoutSec_10day.target
[Unit]
Description=Foo
Description=abc
JobRunningTimeoutSec=10day

