# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_AssertArchitecture_abc.target
[Unit]
Description=Foo
Description=abc
AssertArchitecture=abc

