# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_Documentation_abc.target
[Unit]
Description=Foo
Description=abc
Documentation=abc

