# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_FailureAction_none.target
[Unit]
Description=Foo
Description=abc
FailureAction=none

