# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_Alias_abc.target
[Unit]
Description=Foo
Description=abc
Alias=abc

