# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_ConditionSecurity_abc.target
[Unit]
Description=Foo
Description=abc
ConditionSecurity=abc

