# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_20w.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=20w

