# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobRunningTimeoutSec_9min.target
[Unit]
Description=Foo
Description=abc
JobRunningTimeoutSec=9min

