# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_RefuseManualStop_0.target
[Unit]
Description=Foo
Description=abc
RefuseManualStop=0

