# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_AssertVirtualization_yes.target
[Unit]
Description=Foo
Description=abc
AssertVirtualization=yes

