# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_28h.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=28h

