# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_19d.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=19d

