# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_AssertCPUs_9999999999999.target
[Unit]
Description=Foo
Description=abc
AssertCPUs=9999999999999

