# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_5hr.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=5hr

