# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobRunningTimeoutSec_3hours.target
[Unit]
Description=Foo
Description=abc
JobRunningTimeoutSec=3hours

