# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_ConditionACPower_on.target
[Unit]
Description=Foo
Description=abc
ConditionACPower=on

