# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_RefuseManualStop_off.target
[Unit]
Description=Foo
Description=abc
RefuseManualStop=off

