# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_LogNamespace_abc.target
[Unit]
Description=Foo
Description=abc
LogNamespace=abc

