# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_StopWhenUnneeded_on.target
[Unit]
Description=Foo
Description=abc
StopWhenUnneeded=on

