# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_ConditionFirstBoot_0.target
[Unit]
Description=Foo
Description=abc
ConditionFirstBoot=0

