# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_OnFailureJobMode_isolate.target
[Unit]
Description=Foo
Description=abc
OnFailureJobMode=isolate

