# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_15ms.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=15ms

