# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_ConditionACPower_yes.target
[Unit]
Description=Foo
Description=abc
ConditionACPower=yes

