# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_StartLimitAction_poweroff.target
[Unit]
Description=Foo
Description=abc
StartLimitAction=poweroff

