# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobRunningTimeoutSec_28h.target
[Unit]
Description=Foo
Description=abc
JobRunningTimeoutSec=28h

