# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_26months.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=26months

