# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_17second.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=17second

