# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_12month.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=12month

