# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_ConditionKernelCommandLine_abc.target
[Unit]
Description=Foo
Description=abc
ConditionKernelCommandLine=abc

