# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutRebootArgument_exit-force.target
[Unit]
Description=Foo
Description=abc
JobTimeoutRebootArgument=exit-force

