# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_1%C2%B5s.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=1µs

