# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutAction_exit-force.target
[Unit]
Description=Foo
Description=abc
JobTimeoutAction=exit-force

