# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_AssertFirstBoot_off.target
[Unit]
Description=Foo
Description=abc
AssertFirstBoot=off

