# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_9month.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=9month

