# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_StartLimitBurst_0.target
[Unit]
Description=Foo
Description=abc
StartLimitBurst=0

