# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobRunningTimeoutSec_11y.target
[Unit]
Description=Foo
Description=abc
JobRunningTimeoutSec=11y

