# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_IgnoreOnIsolate_1.target
[Unit]
Description=Foo
Description=abc
IgnoreOnIsolate=1

