# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobRunningTimeoutSec_23m.target
[Unit]
Description=Foo
Description=abc
JobRunningTimeoutSec=23m

