# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobRunningTimeoutSec_4ms.target
[Unit]
Description=Foo
Description=abc
JobRunningTimeoutSec=4ms

