# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_AssertVirtualization_1.target
[Unit]
Description=Foo
Description=abc
AssertVirtualization=1

