# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobRunningTimeoutSec_6days.target
[Unit]
Description=Foo
Description=abc
JobRunningTimeoutSec=6days

