# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_AssertPathIsDirectory_--dev--.target
[Unit]
Description=Foo
Description=abc
AssertPathIsDirectory=/dev/

