# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_ConditionACPower_no.target
[Unit]
Description=Foo
Description=abc
ConditionACPower=no

