# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_AssertFileIsExecutable_--dev--.target
[Unit]
Description=Foo
Description=abc
AssertFileIsExecutable=/dev/

