# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_27msec.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=27msec

