# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_4hour.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=4hour

