# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobRunningTimeoutSec_7seconds.target
[Unit]
Description=Foo
Description=abc
JobRunningTimeoutSec=7seconds

