# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_13minute.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=13minute

