# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_8weeks.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=8weeks

