# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_Description_abc.target
[Unit]
Description=Foo
Description=abc
Description=abc

