# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_AssertKernelVersion_abc.target
[Unit]
Description=Foo
Description=abc
AssertKernelVersion=abc

