# systemdlint/tests/auto/InvalidValue/2.31/good/Unit_StopWhenUnneeded_false.target
[Unit]
Description=Foo
Description=abc
StopWhenUnneeded=false

