# systemdlint/tests/auto/InvalidValue/2.31/good/VLAN_LooseBinding_true.target
[Unit]
Description=Foo
[VLAN]
LooseBinding=true

