# systemdlint/tests/auto/InvalidValue/2.31/good/VLAN_ReorderHeader_off.target
[Unit]
Description=Foo
[VLAN]
ReorderHeader=off

