# systemdlint/tests/auto/InvalidValue/2.31/good/VLAN_ReorderHeader_false.target
[Unit]
Description=Foo
[VLAN]
ReorderHeader=false

