# systemdlint/tests/auto/InvalidValue/2.31/good/VLAN_MVRP_true.target
[Unit]
Description=Foo
[VLAN]
MVRP=true

