# systemdlint/tests/auto/InvalidValue/2.31/good/VLAN_LooseBinding_false.target
[Unit]
Description=Foo
[VLAN]
LooseBinding=false

