# systemdlint/tests/auto/InvalidValue/2.31/good/VLAN_GVRP_true.target
[Unit]
Description=Foo
[VLAN]
GVRP=true

