# systemdlint/tests/auto/InvalidValue/2.31/good/VLAN_ReorderHeader_1.target
[Unit]
Description=Foo
[VLAN]
ReorderHeader=1

