# systemdlint/tests/auto/InvalidValue/2.31/good/VLAN_LooseBinding_off.target
[Unit]
Description=Foo
[VLAN]
LooseBinding=off

