# systemdlint/tests/auto/InvalidValue/2.31/good/VLAN_MVRP_off.target
[Unit]
Description=Foo
[VLAN]
MVRP=off

