# systemdlint/tests/auto/InvalidValue/2.31/good/VLAN_ReorderHeader_yes.target
[Unit]
Description=Foo
[VLAN]
ReorderHeader=yes

