# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_UDPChecksum_0.target
[Unit]
Description=Foo
[VXLAN]
UDPChecksum=0

