# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_ReduceARPProxy_yes.target
[Unit]
Description=Foo
[VXLAN]
ReduceARPProxy=yes

