# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_TOS_abc.target
[Unit]
Description=Foo
[VXLAN]
TOS=abc

