# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_UDPChecksum_true.target
[Unit]
Description=Foo
[VXLAN]
UDPChecksum=true

