# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_RemoteChecksumRx_no.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumRx=no

