# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_Remote_1.1.1.1.target
[Unit]
Description=Foo
[VXLAN]
Remote=1.1.1.1

