# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_ReduceARPProxy_1.target
[Unit]
Description=Foo
[VXLAN]
ReduceARPProxy=1

