# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_GroupPolicyExtension_true.target
[Unit]
Description=Foo
[VXLAN]
GroupPolicyExtension=true

