# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_RouteShortCircuit_0.target
[Unit]
Description=Foo
[VXLAN]
RouteShortCircuit=0

