# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_FDBAgeingSec_0.target
[Unit]
Description=Foo
[VXLAN]
FDBAgeingSec=0

