# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_UDPChecksum_yes.target
[Unit]
Description=Foo
[VXLAN]
UDPChecksum=yes

