# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_UDP6ZeroChecksumRx_no.target
[Unit]
Description=Foo
[VXLAN]
UDP6ZeroChecksumRx=no

