# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_ReduceARPProxy_0.target
[Unit]
Description=Foo
[VXLAN]
ReduceARPProxy=0

