# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_RemoteChecksumTx_off.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumTx=off

