# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_L2MissNotification_false.target
[Unit]
Description=Foo
[VXLAN]
L2MissNotification=false

