# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_RemoteChecksumTx_no.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumTx=no

