# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_UDP6ZeroChecksumTx_no.target
[Unit]
Description=Foo
[VXLAN]
UDP6ZeroChecksumTx=no

