# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_RemoteChecksumRx_false.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumRx=false

