# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_RemoteChecksumRx_yes.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumRx=yes

