# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_FDBAgeingSec_9999999999999.target
[Unit]
Description=Foo
[VXLAN]
FDBAgeingSec=9999999999999

