# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_ReduceARPProxy_no.target
[Unit]
Description=Foo
[VXLAN]
ReduceARPProxy=no

