# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_TTL_255.target
[Unit]
Description=Foo
[VXLAN]
TTL=255

