# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_UDPChecksum_off.target
[Unit]
Description=Foo
[VXLAN]
UDPChecksum=off

