# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_UDP6ZeroChecksumRx_yes.target
[Unit]
Description=Foo
[VXLAN]
UDP6ZeroChecksumRx=yes

