# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_RemoteChecksumRx_off.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumRx=off

