# systemdlint/tests/auto/InvalidValue/2.31/good/VXLAN_ReduceARPProxy_false.target
[Unit]
Description=Foo
[VXLAN]
ReduceARPProxy=false

