# systemdlint/tests/auto/InvalidValue/2.31/good/WireGuard_PrivateKeyFile_--dev--.target
[Unit]
Description=Foo
[WireGuard]
PrivateKeyFile=/dev/

