# systemdlint/tests/auto/InvalidValue/2.31/good/WireGuard_FwMark_abc.target
[Unit]
Description=Foo
[WireGuard]
FwMark=abc

