# systemdlint/tests/auto/InvalidValue/2.31/good/WireGuardPeer_PresharedKeyFile_--dev--.target
[Unit]
Description=Foo
[WireGuardPeer]
PresharedKeyFile=/dev/

