# systemdlint/tests/auto/InvalidValue/2.31/good/WireGuard_ListenPort_65534.target
[Unit]
Description=Foo
[WireGuard]
ListenPort=65534

