# systemdlint/tests/auto/InvalidValue/2.31/good/WireGuard_PrivateKeyFile_--tmp.target
[Unit]
Description=Foo
[WireGuard]
PrivateKeyFile=/tmp

