# systemdlint/tests/auto/InvalidValue/2.31/good/WireGuard_ListenPort_1.target
[Unit]
Description=Foo
[WireGuard]
ListenPort=1

