# tests/auto/OptionTooNew/2.33/bad/Service_StandardInputData_abc.target
[Unit]
Description=Foo
[Service]
StandardInputData=123

