# tests/auto/OptionTooNew/2.33/bad/Service_StandardInput_null.target
[Unit]
Description=Foo
[Service]
StandardInput=123

