# tests/auto/OptionTooNew/2.33/bad/Socket_StandardInputData_abc.target
[Unit]
Description=Foo
[Socket]
StandardInputData=123

