# tests/auto/OptionTooNew/2.33/bad/Swap_LogLevelMax_emerg.target
[Unit]
Description=Foo
[Swap]
LogLevelMax=123

