# tests/auto/OptionTooNew/2.33/bad/Swap_StandardInputData_abc.target
[Unit]
Description=Foo
[Swap]
StandardInputData=123

