# tests/auto/OptionTooNew/2.26/bad/Manager_DefaultTasksAccounting_yes.target
[Unit]
Description=Foo
[Manager]
DefaultTasksAccounting=123

