# tests/auto/OptionTooNew/2.26/bad/Service_TasksAccounting_yes.target
[Unit]
Description=Foo
[Service]
TasksAccounting=123

