# tests/auto/OptionTooNew/2.26/bad/Socket_TasksMax_1.target
[Unit]
Description=Foo
[Socket]
TasksMax=123

