# tests/auto/OptionTooNew/2.26/bad/Socket_Writable_yes.target
[Unit]
Description=Foo
[Socket]
Writable=123

