# tests/auto/OptionTooNew/2.26/bad/Swap_TasksAccounting_yes.target
[Unit]
Description=Foo
[Swap]
TasksAccounting=123

