# systemdlint/tests/auto/InvalidValue/2.47/bad/BridgeMDB_MulticastGroupAddress_255.256.255.255.target
[Unit]
Description=Foo
[BridgeMDB]
MulticastGroupAddress=255.256.255.255

