# systemdlint/tests/auto/InvalidValue/2.47/bad/BridgeMDB_MulticastGroupAddress_True.target
[Unit]
Description=Foo
[BridgeMDB]
MulticastGroupAddress=True

