# systemdlint/tests/auto/InvalidValue/2.47/bad/BridgeMDB_MulticastGroupAddress_1.target
[Unit]
Description=Foo
[BridgeMDB]
MulticastGroupAddress=1

