# systemdlint/tests/auto/InvalidValue/2.47/bad/BridgeMDB_VLANId_-1.target
[Unit]
Description=Foo
[BridgeMDB]
VLANId=-1

