# systemdlint/tests/auto/InvalidValue/2.47/bad/DHCPv6PrefixDelegation_Announce_bar.target
[Unit]
Description=Foo
[DHCPv6PrefixDelegation]
Announce=bar

