# systemdlint/tests/auto/InvalidValue/2.47/bad/FlowQueuePIE_PacketLimit_-1.target
[Unit]
Description=Foo
[FlowQueuePIE]
PacketLimit=-1

