# systemdlint/tests/auto/InvalidValue/2.47/bad/FlowQueuePIE_Handle_0x0.target
[Unit]
Description=Foo
[FlowQueuePIE]
Handle=0x0

