# systemdlint/tests/auto/InvalidValue/2.47/bad/IPv6SendRA_DNSLifetimeSec_0.5d.target
[Unit]
Description=Foo
[IPv6SendRA]
DNSLifetimeSec=0.5d

