# systemdlint/tests/auto/InvalidValue/2.47/bad/Socket_FlushPending_bar.target
[Unit]
Description=Foo
[Socket]
FlushPending=bar

