# systemdlint/tests/auto/InvalidValue/2.47/bad/Timer_FixedRandomDelay_bar.target
[Unit]
Description=Foo
[Timer]
FixedRandomDelay=bar

