# systemdlint/tests/auto/InvalidValue/2.47/bad/Timer_FixedRandomDelay_foo.target
[Unit]
Description=Foo
[Timer]
FixedRandomDelay=foo

