# systemdlint/tests/auto/InvalidValue/2.47/bad/VXLAN_Independent_bar.target
[Unit]
Description=Foo
[VXLAN]
Independent=bar

