# systemdlint/tests/auto/InvalidValue/2.47/bad/VXLAN_Independent_foo.target
[Unit]
Description=Foo
[VXLAN]
Independent=foo

