# systemdlint/tests/auto/OptionTooNew/2.46/good/BFIFO_Handle_0x1.target
[Unit]
Description=Foo
[BFIFO]
Handle=123

