# systemdlint/tests/auto/OptionTooNew/2.46/good/Bridge_VLANProtocol_802.1Q.target
[Unit]
Description=Foo
[Bridge]
VLANProtocol=123

