# systemdlint/tests/auto/OptionTooNew/2.46/good/CAKE_Handle_0x1.target
[Unit]
Description=Foo
[CAKE]
Handle=123

