# systemdlint/tests/auto/OptionTooNew/2.46/good/CAN_SamplePoint_abc.target
[Unit]
Description=Foo
[CAN]
SamplePoint=123

