# systemdlint/tests/auto/OptionTooNew/2.46/good/DeficitRoundRobinSchedulerClass_Parent_abc.target
[Unit]
Description=Foo
[DeficitRoundRobinSchedulerClass]
Parent=123

