# systemdlint/tests/auto/OptionTooNew/2.46/good/DHCPv4_FallbackLeaseLifetimeSec_0.target
[Unit]
Description=Foo
[DHCPv4]
FallbackLeaseLifetimeSec=123

