# systemdlint/tests/auto/OptionTooNew/2.46/good/DHCPv4_SendVendorOption_abc.target
[Unit]
Description=Foo
[DHCPv4]
SendVendorOption=123

