# systemdlint/tests/auto/OptionTooNew/2.46/good/DHCPPrefixDelegation_SubnetId_0x0.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
SubnetId=123

