# systemdlint/tests/auto/OptionTooNew/2.46/good/EnhancedTransmissionSelection_Handle_0x1.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
Handle=123

