# systemdlint/tests/auto/OptionTooNew/2.46/good/GenericRandomEarlyDetection_Handle_0x1.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
Handle=123

