# systemdlint/tests/auto/OptionTooNew/2.46/good/HeavyHitterFilter_Handle_0x1.target
[Unit]
Description=Foo
[HeavyHitterFilter]
Handle=123

