# systemdlint/tests/auto/OptionTooNew/2.46/good/HierarchyTokenBucket_Handle_0x1.target
[Unit]
Description=Foo
[HierarchyTokenBucket]
Handle=123

