# systemdlint/tests/auto/OptionTooNew/2.46/good/Link_TxFlowControl_yes.target
[Unit]
Description=Foo
[Link]
MACAddress=abc
TxFlowControl=123

