# systemdlint/tests/auto/OptionTooNew/2.46/good/LLDP_MUDURL_abc.target
[Unit]
Description=Foo
[LLDP]
MUDURL=123

