# systemdlint/tests/auto/OptionTooNew/2.46/good/PFIFOFast_Handle_0x1.target
[Unit]
Description=Foo
[PFIFOFast]
Handle=123

