# systemdlint/tests/auto/OptionTooNew/2.46/good/PFIFO_Handle_0x1.target
[Unit]
Description=Foo
[PFIFO]
Handle=123

