# systemdlint/tests/auto/OptionTooNew/2.46/good/Socket_RootHash_abc.target
[Unit]
Description=Foo
[Socket]
RootHash=123

