# systemdlint/tests/auto/OptionTooNew/2.46/good/SR-IOV_VLANId_1.target
[Unit]
Description=Foo
[SR-IOV]
VLANId=123

