# systemdlint/tests/auto/OptionTooNew/2.46/good/StochasticFairBlue_Parent_abc.target
[Unit]
Description=Foo
[StochasticFairBlue]
Parent=123

