# tests/manual/Multiplicity/2.31/good/Service_Simple_ExecStart.service
[Unit]
Description=one ExecStart

[Service]
Type=simple
ExecStart=/bin/busybox