# tests/manual/Multiplicity/2.31/good/Service_OneShot_ExecStart.service
[Unit]
Description=at least one ExecStart

[Service]
Type=oneshot
ExecStart=/bin/busybox