# systemdlint/tests/auto/OptionTooNew/2.34/bad/Mount_StateDirectory_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
StateDirectory=123

