# systemdlint/tests/auto/OptionTooNew/2.34/bad/Service_IPAddressAllow_1.1.1.1--0.target
[Unit]
Description=Foo
[Service]
IPAddressAllow=123

