# systemdlint/tests/auto/OptionTooNew/2.34/bad/Service_LogsDirectory_abc.target
[Unit]
Description=Foo
[Service]
LogsDirectory=123

