# systemdlint/tests/auto/OptionTooNew/2.34/bad/Service_KeyringMode_inherit.target
[Unit]
Description=Foo
[Service]
KeyringMode=123

