# systemdlint/tests/auto/OptionTooNew/2.34/bad/Service_ConfigurationDirectory_abc.target
[Unit]
Description=Foo
[Service]
ConfigurationDirectory=123

