# systemdlint/tests/auto/OptionTooNew/2.34/bad/Socket_RuntimeDirectory_abc.target
[Unit]
Description=Foo
[Socket]
RuntimeDirectory=123

