# systemdlint/tests/auto/OptionTooNew/2.34/bad/Socket_ConfigurationDirectory_abc.target
[Unit]
Description=Foo
[Socket]
ConfigurationDirectory=123

