# systemdlint/tests/auto/OptionTooNew/2.34/bad/Swap_ConfigurationDirectoryMode_0.target
[Unit]
Description=Foo
[Swap]
ConfigurationDirectoryMode=123

