# systemdlint/tests/auto/OptionTooNew/2.34/bad/Swap_StateDirectoryMode_0.target
[Unit]
Description=Foo
[Swap]
StateDirectoryMode=123

