# systemdlint/tests/auto/OptionTooNew/2.34/bad/Swap_IPAddressAllow_1.1.1.1--0.target
[Unit]
Description=Foo
[Swap]
IPAddressAllow=123

