# tests/auto/InvalidValue/2.35/good/Mount_IPAddressDeny_1.1.1.1--0.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
IPAddressDeny=1.1.1.1/0

