# tests/auto/InvalidValue/2.35/good/Mount_StateDirectory_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
StateDirectory=abc

