# tests/auto/InvalidValue/2.35/good/Service_RuntimeDirectoryPreserve_1.target
[Unit]
Description=Foo
[Service]
RuntimeDirectoryPreserve=1

