# tests/auto/InvalidValue/2.35/good/Service_ConfigurationDirectory_abc.target
[Unit]
Description=Foo
[Service]
ConfigurationDirectory=abc

