# tests/auto/InvalidValue/2.35/good/Service_RuntimeDirectoryPreserve_true.target
[Unit]
Description=Foo
[Service]
RuntimeDirectoryPreserve=true

