# tests/auto/InvalidValue/2.35/good/Service_IPAccounting_off.target
[Unit]
Description=Foo
[Service]
IPAccounting=off

