# tests/auto/InvalidValue/2.35/good/Service_LogsDirectory_abc.target
[Unit]
Description=Foo
[Service]
LogsDirectory=abc

