# tests/auto/InvalidValue/2.35/good/Service_RuntimeDirectoryPreserve_yes.target
[Unit]
Description=Foo
[Service]
RuntimeDirectoryPreserve=yes

