# tests/auto/InvalidValue/2.35/good/Service_LockPersonality_on.target
[Unit]
Description=Foo
[Service]
LockPersonality=on

