# tests/auto/InvalidValue/2.35/good/Service_CacheDirectoryMode_9999999999999.target
[Unit]
Description=Foo
[Service]
CacheDirectoryMode=9999999999999

