# tests/auto/InvalidValue/2.35/good/Service_IPAccounting_true.target
[Unit]
Description=Foo
[Service]
IPAccounting=true

