# tests/auto/InvalidValue/2.35/good/Service_LockPersonality_true.target
[Unit]
Description=Foo
[Service]
LockPersonality=true

