# tests/auto/InvalidValue/2.35/good/Service_KeyringMode_shared.target
[Unit]
Description=Foo
[Service]
KeyringMode=shared

