# tests/auto/InvalidValue/2.35/good/Service_IPAccounting_on.target
[Unit]
Description=Foo
[Service]
IPAccounting=on

