# tests/auto/InvalidValue/2.35/good/Service_StateDirectory_abc.target
[Unit]
Description=Foo
[Service]
StateDirectory=abc

