# tests/auto/InvalidValue/2.35/good/Service_CacheDirectory_abc.target
[Unit]
Description=Foo
[Service]
CacheDirectory=abc

