# tests/auto/InvalidValue/2.35/good/Service_KeyringMode_private.target
[Unit]
Description=Foo
[Service]
KeyringMode=private

