# tests/auto/InvalidValue/2.35/good/Service_RuntimeDirectoryPreserve_no.target
[Unit]
Description=Foo
[Service]
RuntimeDirectoryPreserve=no

