# tests/auto/InvalidValue/2.35/good/Service_LockPersonality_false.target
[Unit]
Description=Foo
[Service]
LockPersonality=false

