# tests/auto/InvalidValue/2.35/good/Service_IPAccounting_yes.target
[Unit]
Description=Foo
[Service]
IPAccounting=yes

