# tests/auto/InvalidValue/2.35/good/Service_IPAccounting_no.target
[Unit]
Description=Foo
[Service]
IPAccounting=no

