# tests/auto/InvalidValue/2.35/good/Socket_RuntimeDirectoryMode_9999999999999.target
[Unit]
Description=Foo
[Socket]
RuntimeDirectoryMode=9999999999999

