# tests/auto/InvalidValue/2.35/good/Socket_KeyringMode_private.target
[Unit]
Description=Foo
[Socket]
KeyringMode=private

