# tests/auto/InvalidValue/2.35/good/Socket_RuntimeDirectoryPreserve_true.target
[Unit]
Description=Foo
[Socket]
RuntimeDirectoryPreserve=true

