# tests/auto/InvalidValue/2.35/good/Socket_LockPersonality_on.target
[Unit]
Description=Foo
[Socket]
LockPersonality=on

