# tests/auto/InvalidValue/2.35/good/Socket_IPAccounting_on.target
[Unit]
Description=Foo
[Socket]
IPAccounting=on

