# tests/auto/InvalidValue/2.35/good/Socket_KeyringMode_inherit.target
[Unit]
Description=Foo
[Socket]
KeyringMode=inherit

