# tests/auto/InvalidValue/2.35/good/Socket_LockPersonality_no.target
[Unit]
Description=Foo
[Socket]
LockPersonality=no

