# tests/auto/InvalidValue/2.35/good/Socket_RuntimeDirectoryPreserve_yes.target
[Unit]
Description=Foo
[Socket]
RuntimeDirectoryPreserve=yes

