# tests/auto/InvalidValue/2.35/good/Socket_StateDirectoryMode_9999999999999.target
[Unit]
Description=Foo
[Socket]
StateDirectoryMode=9999999999999

