# tests/auto/InvalidValue/2.35/good/Socket_RuntimeDirectoryPreserve_false.target
[Unit]
Description=Foo
[Socket]
RuntimeDirectoryPreserve=false

