# tests/auto/InvalidValue/2.35/good/Socket_RuntimeDirectoryPreserve_off.target
[Unit]
Description=Foo
[Socket]
RuntimeDirectoryPreserve=off

