# tests/auto/InvalidValue/2.35/good/Socket_IPAccounting_true.target
[Unit]
Description=Foo
[Socket]
IPAccounting=true

