# tests/auto/InvalidValue/2.35/good/Socket_RuntimeDirectoryPreserve_0.target
[Unit]
Description=Foo
[Socket]
RuntimeDirectoryPreserve=0

