# tests/auto/InvalidValue/2.35/good/Socket_RuntimeDirectory_abc.target
[Unit]
Description=Foo
[Socket]
RuntimeDirectory=abc

