# tests/auto/InvalidValue/2.35/good/Swap_ConfigurationDirectoryMode_9999999999999.target
[Unit]
Description=Foo
[Swap]
ConfigurationDirectoryMode=9999999999999

