# tests/auto/InvalidValue/2.35/good/Swap_IPAccounting_no.target
[Unit]
Description=Foo
[Swap]
IPAccounting=no

