# tests/auto/InvalidValue/2.35/good/Swap_StateDirectoryMode_9999999999999.target
[Unit]
Description=Foo
[Swap]
StateDirectoryMode=9999999999999

