# tests/auto/InvalidValue/2.35/good/Swap_LockPersonality_no.target
[Unit]
Description=Foo
[Swap]
LockPersonality=no

