# tests/auto/InvalidValue/2.35/good/Swap_UnsetEnvironment_abc.target
[Unit]
Description=Foo
[Swap]
UnsetEnvironment=abc

