# tests/auto/InvalidValue/2.35/good/Swap_LockPersonality_false.target
[Unit]
Description=Foo
[Swap]
LockPersonality=false

