# tests/auto/InvalidValue/2.35/good/Swap_LogsDirectory_abc.target
[Unit]
Description=Foo
[Swap]
LogsDirectory=abc

