# tests/auto/InvalidValue/2.35/good/Swap_LockPersonality_yes.target
[Unit]
Description=Foo
[Swap]
LockPersonality=yes

