# tests/auto/InvalidValue/2.35/good/Swap_RuntimeDirectory_abc.target
[Unit]
Description=Foo
[Swap]
RuntimeDirectory=abc

