# tests/auto/InvalidValue/2.35/good/Swap_LockPersonality_1.target
[Unit]
Description=Foo
[Swap]
LockPersonality=1

