# tests/auto/InvalidValue/2.35/good/Swap_RuntimeDirectoryPreserve_on.target
[Unit]
Description=Foo
[Swap]
RuntimeDirectoryPreserve=on

