# tests/auto/InvalidValue/2.35/good/Swap_KeyringMode_shared.target
[Unit]
Description=Foo
[Swap]
KeyringMode=shared

