# tests/auto/InvalidValue/2.35/good/Swap_RuntimeDirectoryMode_9999999999999.target
[Unit]
Description=Foo
[Swap]
RuntimeDirectoryMode=9999999999999

