# tests/auto/InvalidValue/2.35/good/Swap_ConfigurationDirectoryMode_0.target
[Unit]
Description=Foo
[Swap]
ConfigurationDirectoryMode=0

