# tests/auto/InvalidValue/2.35/good/Swap_IPAddressAllow_255.255.255.255--32.target
[Unit]
Description=Foo
[Swap]
IPAddressAllow=255.255.255.255/32

