# tests/auto/InvalidValue/2.35/good/Swap_CacheDirectoryMode_9999999999999.target
[Unit]
Description=Foo
[Swap]
CacheDirectoryMode=9999999999999

