# tests/auto/InvalidValue/2.35/good/Swap_IPAccounting_yes.target
[Unit]
Description=Foo
[Swap]
IPAccounting=yes

