# tests/auto/InvalidValue/2.35/good/Swap_LockPersonality_off.target
[Unit]
Description=Foo
[Swap]
LockPersonality=off

