# tests/auto/InvalidValue/2.48/good/IPv6AcceptRA_PrefixAllowList_abc.target
[Unit]
Description=Foo
[IPv6AcceptRA]
PrefixAllowList=abc

