# tests/auto/InvalidValue/2.48/good/Link_ActivationPolicy_always-up.target
[Unit]
Description=Foo
[Link]
ActivationPolicy=always-up

