# tests/auto/InvalidValue/2.48/good/Link_GenericSegmentOffloadMaxSegments_1G.target
[Unit]
Description=Foo
[Link]
GenericSegmentOffloadMaxSegments=1G

