# tests/auto/InvalidValue/2.48/good/Link_GenericSegmentOffloadMaxSegments_1M.target
[Unit]
Description=Foo
[Link]
GenericSegmentOffloadMaxSegments=1M

