# tests/auto/InvalidValue/2.48/good/Link_GenericSegmentOffloadMaxBytes_65536M.target
[Unit]
Description=Foo
[Link]
GenericSegmentOffloadMaxBytes=65536M

