# tests/auto/InvalidValue/2.48/good/Link_GenericSegmentOffloadMaxSegments_65535G.target
[Unit]
Description=Foo
[Link]
GenericSegmentOffloadMaxSegments=65535G

