# tests/auto/InvalidValue/2.48/good/Link_ActivationPolicy_up.target
[Unit]
Description=Foo
[Link]
ActivationPolicy=up

