# tests/auto/InvalidValue/2.48/good/Link_TransmitQueues_4096.target
[Unit]
Description=Foo
[Link]
TransmitQueues=4096

