# tests/auto/InvalidValue/2.48/good/Link_GenericSegmentOffloadMaxBytes_1M.target
[Unit]
Description=Foo
[Link]
GenericSegmentOffloadMaxBytes=1M

