# tests/auto/InvalidValue/2.48/good/Network_RouteTable_abc.target
[Unit]
Description=Foo
[Network]
RouteTable=abc

