# tests/auto/InvalidValue/2.48/good/OOM_DefaultMemoryPressureDurationSec_9999999999999.target
[Unit]
Description=Foo
[OOM]
DefaultMemoryPressureDurationSec=9999999999999

