# tests/auto/InvalidValue/2.48/good/RoutingPolicyRule_Type_unreachable.target
[Unit]
Description=Foo
[RoutingPolicyRule]
Type=unreachable

