# tests/auto/InvalidValue/2.48/good/Scope_ManagedOOMPreference_none.target
[Unit]
Description=Foo
[Scope]
ManagedOOMPreference=none

