# tests/auto/InvalidValue/2.48/good/Scope_ManagedOOMPreference_omit.target
[Unit]
Description=Foo
[Scope]
ManagedOOMPreference=omit

