# tests/auto/InvalidValue/2.48/good/Service_ExecPaths_abc.target
[Unit]
Description=Foo
[Service]
ExecPaths=abc

