# tests/auto/InvalidValue/2.48/good/Service_ManagedOOMPreference_none.target
[Unit]
Description=Foo
[Service]
ManagedOOMPreference=none

