# tests/auto/InvalidValue/2.48/good/Service_ManagedOOMPreference_avoid.target
[Unit]
Description=Foo
[Service]
ManagedOOMPreference=avoid

