# tests/auto/InvalidValue/2.48/good/Slice_ManagedOOMPreference_none.target
[Unit]
Description=Foo
[Slice]
ManagedOOMPreference=none

