# tests/auto/InvalidValue/2.48/good/Slice_ManagedOOMPreference_avoid.target
[Unit]
Description=Foo
[Slice]
ManagedOOMPreference=avoid

