# tests/auto/InvalidValue/2.48/good/Socket_NoExecPaths_abc.target
[Unit]
Description=Foo
[Socket]
NoExecPaths=abc

