# tests/auto/InvalidValue/2.48/good/Socket_ManagedOOMPreference_omit.target
[Unit]
Description=Foo
[Socket]
ManagedOOMPreference=omit

