# tests/auto/InvalidValue/2.48/good/Swap_NoExecPaths_abc.target
[Unit]
Description=Foo
[Swap]
NoExecPaths=abc

