# tests/auto/InvalidValue/2.48/good/Swap_ExecPaths_abc.target
[Unit]
Description=Foo
[Swap]
ExecPaths=abc

