# tests/auto/InvalidValue/2.48/good/Swap_ManagedOOMPreference_omit.target
[Unit]
Description=Foo
[Swap]
ManagedOOMPreference=omit

