# tests/auto/InvalidValue/2.48/good/Time_ConnectionRetrySec_1.target
[Unit]
Description=Foo
[Time]
ConnectionRetrySec=1

