# tests/auto/InvalidValue/2.48/good/Unit_ConditionCPUFeature_rdseed.target
[Unit]
Description=Foo
Description=abc
ConditionCPUFeature=rdseed

