# tests/auto/InvalidValue/2.48/good/Unit_ConditionCPUFeature_rdtscp.target
[Unit]
Description=Foo
Description=abc
ConditionCPUFeature=rdtscp

