# tests/auto/InvalidValue/2.48/good/Unit_ConditionCPUFeature_lahf_lm.target
[Unit]
Description=Foo
Description=abc
ConditionCPUFeature=lahf_lm

