# tests/auto/InvalidValue/2.48/good/Unit_ConditionCPUFeature_bmi2.target
[Unit]
Description=Foo
Description=abc
ConditionCPUFeature=bmi2

