# tests/auto/InvalidValue/2.48/good/Unit_ConditionCPUFeature_sha_ni.target
[Unit]
Description=Foo
Description=abc
ConditionCPUFeature=sha_ni

