# tests/auto/InvalidValue/2.48/good/VLAN_EgressQOSMaps_abc.target
[Unit]
Description=Foo
[VLAN]
EgressQOSMaps=abc

