# tests/auto/InvalidValue/2.48/good/VLAN_IngressQOSMaps_abc.target
[Unit]
Description=Foo
[VLAN]
IngressQOSMaps=abc

