# systemdlint/tests/auto/OptionDeprecated/2.28/bad/Manager_CrashChVT_abc.target
[Unit]
Description=Foo
[Manager]
CrashChVT=123

