# systemdlint/tests/auto/OptionTooNew/2.39/bad/Service_MemoryMin_0K.target
[Unit]
Description=Foo
[Service]
MemoryMin=123

