# systemdlint/tests/auto/OptionTooNew/2.09/bad/Socket_AppArmorProfile_abc.target
[Unit]
Description=Foo
[Socket]
AppArmorProfile=123

