# tests/auto/InvalidValue/2.47/good/BridgeMDB_MulticastGroupAddress_255.255.255.255--32.target
[Unit]
Description=Foo
[BridgeMDB]
MulticastGroupAddress=255.255.255.255/32

