# tests/auto/InvalidValue/2.47/good/BridgeMDB_MulticastGroupAddress_file%3A----test.a.b.com.target
[Unit]
Description=Foo
[BridgeMDB]
MulticastGroupAddress=file://test.a.b.com

