# tests/auto/InvalidValue/2.47/good/BridgeMDB_MulticastGroupAddress_http%3A----user%3Apwd%40test.com%3A1234.target
[Unit]
Description=Foo
[BridgeMDB]
MulticastGroupAddress=http://user:pwd@test.com:1234

