# tests/auto/InvalidValue/2.47/good/BridgeMDB_VLANId_4094.target
[Unit]
Description=Foo
[BridgeMDB]
VLANId=4094

