# tests/auto/InvalidValue/2.47/good/DHCPv6PrefixDelegation_Announce_0.target
[Unit]
Description=Foo
[DHCPv6PrefixDelegation]
Announce=0

