# tests/auto/InvalidValue/2.47/good/FlowQueuePIE_PacketLimit_0.target
[Unit]
Description=Foo
[FlowQueuePIE]
PacketLimit=0

