# tests/auto/InvalidValue/2.47/good/FlowQueuePIE_Parent_root.target
[Unit]
Description=Foo
[FlowQueuePIE]
Parent=root

