# tests/auto/InvalidValue/2.47/good/FlowQueuePIE_Parent_ingress.target
[Unit]
Description=Foo
[FlowQueuePIE]
Parent=ingress

