# tests/auto/InvalidValue/2.47/good/FlowQueuePIE_Handle_0x1.target
[Unit]
Description=Foo
[FlowQueuePIE]
Handle=0x1

