# tests/auto/InvalidValue/2.47/good/IPv6SendRA_DNSLifetimeSec_26year.target
[Unit]
Description=Foo
[IPv6SendRA]
DNSLifetimeSec=26year

