# tests/auto/InvalidValue/2.47/good/IPv6SendRA_EmitDomains_off.target
[Unit]
Description=Foo
[IPv6SendRA]
EmitDomains=off

