# tests/auto/InvalidValue/2.47/good/IPv6SendRA_OtherInformation_1.target
[Unit]
Description=Foo
[IPv6SendRA]
OtherInformation=1

