# tests/auto/InvalidValue/2.47/good/IPv6SendRA_EmitDNS_on.target
[Unit]
Description=Foo
[IPv6SendRA]
EmitDNS=on

