# tests/auto/InvalidValue/2.47/good/IPv6SendRA_DNSLifetimeSec_17second.target
[Unit]
Description=Foo
[IPv6SendRA]
DNSLifetimeSec=17second

