# tests/auto/InvalidValue/2.47/good/IPv6SendRA_DNSLifetimeSec_10M.target
[Unit]
Description=Foo
[IPv6SendRA]
DNSLifetimeSec=10M

