# tests/auto/InvalidValue/2.47/good/IPv6SendRA_DNSLifetimeSec_24w.target
[Unit]
Description=Foo
[IPv6SendRA]
DNSLifetimeSec=24w

