# tests/auto/InvalidValue/2.47/good/IPv6SendRA_DNSLifetimeSec_5hr.target
[Unit]
Description=Foo
[IPv6SendRA]
DNSLifetimeSec=5hr

