# tests/auto/InvalidValue/2.47/good/IPv6SendRA_Domains_abc.target
[Unit]
Description=Foo
[IPv6SendRA]
Domains=abc

