# tests/auto/InvalidValue/2.47/good/IPv6SendRA_RouterLifetimeSec_5hr.target
[Unit]
Description=Foo
[IPv6SendRA]
RouterLifetimeSec=5hr

