# tests/auto/InvalidValue/2.47/good/IPv6SendRA_DNSLifetimeSec_28%C2%B5s.target
[Unit]
Description=Foo
[IPv6SendRA]
DNSLifetimeSec=28µs

