# tests/auto/InvalidValue/2.47/good/IPv6SendRA_EmitDNS_off.target
[Unit]
Description=Foo
[IPv6SendRA]
EmitDNS=off

