# tests/auto/InvalidValue/2.47/good/IPv6SendRA_RouterLifetimeSec_12month.target
[Unit]
Description=Foo
[IPv6SendRA]
RouterLifetimeSec=12month

