# tests/auto/InvalidValue/2.47/good/IPv6SendRA_RouterPreference_high.target
[Unit]
Description=Foo
[IPv6SendRA]
RouterPreference=high

