# tests/auto/InvalidValue/2.47/good/IPv6SendRA_DNSLifetimeSec_18sec.target
[Unit]
Description=Foo
[IPv6SendRA]
DNSLifetimeSec=18sec

