# tests/auto/InvalidValue/2.47/good/Mount_ManagedOOMSwap_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
ManagedOOMSwap=abc

