# tests/auto/InvalidValue/2.47/good/Network_DHCPv6PrefixDelegation_yes.target
[Unit]
Description=Foo
[Network]
DHCPv6PrefixDelegation=yes

