# tests/auto/InvalidValue/2.47/good/Network_DHCPv6PrefixDelegation_on.target
[Unit]
Description=Foo
[Network]
DHCPv6PrefixDelegation=on

