# tests/auto/InvalidValue/2.47/good/Network_IPv6SendRA_dhcpv6.target
[Unit]
Description=Foo
[Network]
IPv6SendRA=dhcpv6

