# tests/auto/InvalidValue/2.47/good/Network_DHCPv6PrefixDelegation_off.target
[Unit]
Description=Foo
[Network]
DHCPv6PrefixDelegation=off

