# tests/auto/InvalidValue/2.47/good/OOM_SwapUsedLimit_abc.target
[Unit]
Description=Foo
[OOM]
SwapUsedLimit=abc

