# tests/auto/InvalidValue/2.47/good/Service_ProtectProc_noaccess.target
[Unit]
Description=Foo
[Service]
ProtectProc=noaccess

