# tests/auto/InvalidValue/2.47/good/Service_NUMAMask_all.target
[Unit]
Description=Foo
[Service]
NUMAMask=all

