# tests/auto/InvalidValue/2.47/good/Service_ProcSubset_pid.target
[Unit]
Description=Foo
[Service]
ProcSubset=pid

