# tests/auto/InvalidValue/2.47/good/Service_ManagedOOMMemoryPressure_abc.target
[Unit]
Description=Foo
[Service]
ManagedOOMMemoryPressure=abc

