# tests/auto/InvalidValue/2.47/good/Service_RootImageOptions_%20root%20root-secondary%20home%20srv%20esp%20xbootldr%20tmp%20var%20usr.target
[Unit]
Description=Foo
[Service]
RootImageOptions= root root-secondary home srv esp xbootldr tmp var usr

