# tests/auto/InvalidValue/2.47/good/Service_NUMAMask_abc.target
[Unit]
Description=Foo
[Service]
NUMAMask=abc

