# tests/auto/InvalidValue/2.47/good/Service_ProtectProc_ptraceable.target
[Unit]
Description=Foo
[Service]
ProtectProc=ptraceable

