# tests/auto/InvalidValue/2.47/good/Socket_FlushPending_0.target
[Unit]
Description=Foo
[Socket]
FlushPending=0

