# tests/auto/InvalidValue/2.47/good/Socket_ManagedOOMMemoryPressureLimitPercent_abc.target
[Unit]
Description=Foo
[Socket]
ManagedOOMMemoryPressureLimitPercent=abc

