# tests/auto/InvalidValue/2.47/good/Socket_FlushPending_off.target
[Unit]
Description=Foo
[Socket]
FlushPending=off

