# tests/auto/InvalidValue/2.47/good/Socket_ManagedOOMSwap_abc.target
[Unit]
Description=Foo
[Socket]
ManagedOOMSwap=abc

