# tests/auto/InvalidValue/2.47/good/Socket_Timestamping_%C2%B5s.target
[Unit]
Description=Foo
[Socket]
Timestamping=µs

