# tests/auto/InvalidValue/2.47/good/Socket_ProtectProc_ptraceable.target
[Unit]
Description=Foo
[Socket]
ProtectProc=ptraceable

