# tests/auto/InvalidValue/2.47/good/Socket_ProtectProc_noaccess.target
[Unit]
Description=Foo
[Socket]
ProtectProc=noaccess

