# tests/auto/InvalidValue/2.47/good/Socket_FlushPending_on.target
[Unit]
Description=Foo
[Socket]
FlushPending=on

