# tests/auto/InvalidValue/2.47/good/Socket_ProtectProc_invisible.target
[Unit]
Description=Foo
[Socket]
ProtectProc=invisible

