# tests/auto/InvalidValue/2.47/good/Socket_FlushPending_true.target
[Unit]
Description=Foo
[Socket]
FlushPending=true

