# tests/auto/InvalidValue/2.47/good/Socket_ProcSubset_pid.target
[Unit]
Description=Foo
[Socket]
ProcSubset=pid

