# tests/auto/InvalidValue/2.47/good/Socket_ProcSubset_all.target
[Unit]
Description=Foo
[Socket]
ProcSubset=all

