# tests/auto/InvalidValue/2.47/good/Swap_ManagedOOMMemoryPressure_abc.target
[Unit]
Description=Foo
[Swap]
ManagedOOMMemoryPressure=abc

