# tests/auto/InvalidValue/2.47/good/Swap_ProtectProc_default.target
[Unit]
Description=Foo
[Swap]
ProtectProc=default

