# tests/auto/InvalidValue/2.47/good/Swap_NUMAMask_abc.target
[Unit]
Description=Foo
[Swap]
NUMAMask=abc

