# tests/auto/InvalidValue/2.47/good/Swap_NUMAMask_all.target
[Unit]
Description=Foo
[Swap]
NUMAMask=all

