# tests/auto/InvalidValue/2.47/good/Swap_ProtectProc_invisible.target
[Unit]
Description=Foo
[Swap]
ProtectProc=invisible

