# tests/auto/InvalidValue/2.47/good/Swap_MountImages_abc.target
[Unit]
Description=Foo
[Swap]
MountImages=abc

