# tests/auto/InvalidValue/2.47/good/Timer_FixedRandomDelay_off.target
[Unit]
Description=Foo
[Timer]
FixedRandomDelay=off

