# tests/auto/InvalidValue/2.47/good/Timer_FixedRandomDelay_1.target
[Unit]
Description=Foo
[Timer]
FixedRandomDelay=1

