# tests/auto/InvalidValue/2.47/good/VXLAN_Independent_0.target
[Unit]
Description=Foo
[VXLAN]
Independent=0

