# tests/auto/InvalidValue/2.47/good/VXLAN_Independent_on.target
[Unit]
Description=Foo
[VXLAN]
Independent=on

