# tests/auto/InvalidValue/2.47/good/VXLAN_Independent_1.target
[Unit]
Description=Foo
[VXLAN]
Independent=1

