# tests/auto/InvalidValue/2.47/good/VXLAN_Independent_true.target
[Unit]
Description=Foo
[VXLAN]
Independent=true

