# tests/auto/OptionTooNew/2.44/good/DHCPServer_EmitSIP_yes.target
[Unit]
Description=Foo
[DHCPServer]
EmitSIP=123

