# tests/auto/OptionTooNew/2.44/good/DHCPv6_PrefixDelegationHint_1.1.1.1--0.target
[Unit]
Description=Foo
[DHCPv6]
PrefixDelegationHint=123

