# tests/auto/OptionTooNew/2.44/good/DHCPv4_SendOption_abc.target
[Unit]
Description=Foo
[DHCPv4]
SendOption=123

