# tests/auto/OptionTooNew/2.44/good/Route_IPServiceType_CS4.target
[Unit]
Description=Foo
[Route]
IPServiceType=123

