# tests/auto/OptionTooNew/2.44/good/Service_RestartKillSignal_0.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=123

