# tests/auto/OptionTooNew/2.09/bad/Service_AppArmorProfile_abc.target
[Unit]
Description=Foo
[Service]
AppArmorProfile=123

