# tests/auto/OptionTooNew/2.09/bad/Unit_ConditionArchitecture_abc.target
[Unit]
Description=Foo
Description=abc
ConditionArchitecture=123

