# systemdlint/tests/auto/OptionDeprecated/2.47/bad/Network_IPv6PrefixDelegation_static.target
[Unit]
Description=Foo
[Network]
IPv6PrefixDelegation=123

