# systemdlint/tests/auto/MandatoryOptionMissing/2.31/bad/IPv6AddressLabel_Mandatory_IPv6AddressLabel.target
[Unit]
Description=Foo
[IPv6AddressLabel]
