# systemdlint/tests/auto/InvalidValue/2.49/good/DHCPv6PrefixDelegation_RouteMetric_4294967295.target
[Unit]
Description=Foo
[DHCPv6PrefixDelegation]
RouteMetric=4294967295

