# systemdlint/tests/auto/InvalidValue/2.49/good/DHCPServer_RelayTarget_1.1.1.1.target
[Unit]
Description=Foo
[DHCPServer]
RelayTarget=1.1.1.1

