# systemdlint/tests/auto/InvalidValue/2.49/good/DHCPServerStaticLease_MACAddress_AABB.CCDD.EEFF.target
[Unit]
Description=Foo
[DHCPServerStaticLease]
MACAddress=AABB.CCDD.EEFF

