# systemdlint/tests/auto/InvalidValue/2.49/good/DHCPv6PrefixDelegation_RouteMetric_0.target
[Unit]
Description=Foo
[DHCPv6PrefixDelegation]
RouteMetric=0

