# systemdlint/tests/auto/InvalidValue/2.49/good/DHCPServer_ServerAddress_1.1.1.1.target
[Unit]
Description=Foo
[DHCPServer]
ServerAddress=1.1.1.1

