# systemdlint/tests/auto/InvalidValue/2.49/good/IPv6AcceptRA_RouteMetric_0.target
[Unit]
Description=Foo
[IPv6AcceptRA]
RouteMetric=0

