# systemdlint/tests/auto/InvalidValue/2.49/good/Link_RequiredFamilyForOnline_ipv4.target
[Unit]
Description=Foo
[Link]
RequiredFamilyForOnline=ipv4

