# systemdlint/tests/auto/InvalidValue/2.49/good/Link_RequiredFamilyForOnline_both.target
[Unit]
Description=Foo
[Link]
RequiredFamilyForOnline=both

