# systemdlint/tests/auto/InvalidValue/2.49/good/Link_RequiredFamilyForOnline_ipv6.target
[Unit]
Description=Foo
[Link]
RequiredFamilyForOnline=ipv6

