# systemdlint/tests/auto/InvalidValue/2.49/good/Manager_StatusUnitFormat_name.target
[Unit]
Description=Foo
[Manager]
StatusUnitFormat=name

