# systemdlint/tests/auto/InvalidValue/2.49/good/Network_ManageForeignRoutingPolicyRules_yes.target
[Unit]
Description=Foo
[Network]
ManageForeignRoutingPolicyRules=yes

