# systemdlint/tests/auto/InvalidValue/2.49/good/Network_ManageForeignRoutingPolicyRules_on.target
[Unit]
Description=Foo
[Network]
ManageForeignRoutingPolicyRules=on

