# systemdlint/tests/auto/InvalidValue/2.49/good/Network_IPv6StableSecretAddress_abc.target
[Unit]
Description=Foo
[Network]
IPv6StableSecretAddress=abc

