# systemdlint/tests/auto/InvalidValue/2.49/good/Network_ManageForeignRoutingPolicyRules_no.target
[Unit]
Description=Foo
[Network]
ManageForeignRoutingPolicyRules=no

