# systemdlint/tests/auto/InvalidValue/2.49/good/Network_ManageForeignRoutingPolicyRules_false.target
[Unit]
Description=Foo
[Network]
ManageForeignRoutingPolicyRules=false

