# systemdlint/tests/auto/InvalidValue/2.49/good/Network_ManageForeignRoutingPolicyRules_true.target
[Unit]
Description=Foo
[Network]
ManageForeignRoutingPolicyRules=true

