# systemdlint/tests/auto/InvalidValue/2.49/good/Scope_SocketBindDeny_abc.target
[Unit]
Description=Foo
[Scope]
SocketBindDeny=abc

