# systemdlint/tests/auto/InvalidValue/2.49/good/Service_SocketBindDeny_abc.target
[Unit]
Description=Foo
[Service]
SocketBindDeny=abc

