# systemdlint/tests/auto/InvalidValue/2.49/good/Service_SocketBindAllow_abc.target
[Unit]
Description=Foo
[Service]
SocketBindAllow=abc

