# systemdlint/tests/auto/InvalidValue/2.49/good/Unit_BPFProgram_abc.target
[Unit]
Description=Foo
Description=abc
BPFProgram=abc

