# systemdlint/tests/auto/InvalidValue/2.49/good/Upload_NetworkTimeoutSec_13months.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=13months

