# systemdlint/tests/auto/InvalidValue/2.49/good/Upload_NetworkTimeoutSec_6h.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=6h

