# systemdlint/tests/auto/InvalidValue/2.49/good/Upload_NetworkTimeoutSec_17second.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=17second

