# systemdlint/tests/auto/InvalidValue/2.49/good/Upload_NetworkTimeoutSec_21us.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=21us

