# systemdlint/tests/auto/InvalidValue/2.49/good/Upload_NetworkTimeoutSec_8minute.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=8minute

