# systemdlint/tests/auto/InvalidValue/2.49/good/Upload_NetworkTimeoutSec_18sec.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=18sec

