# systemdlint/tests/auto/InvalidValue/2.49/good/Upload_URL_file%3A----test.a.b.com.target
[Unit]
Description=Foo
[Upload]
URL=file://test.a.b.com

