# systemdlint/tests/auto/InvalidValue/2.49/good/Upload_ServerKeyFile_abc.target
[Unit]
Description=Foo
[Upload]
ServerKeyFile=abc

