# systemdlint/tests/auto/InvalidValue/2.49/good/Upload_NetworkTimeoutSec_23week.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=23week

