# tests/auto/InvalidValue/2.48/bad/Link_GenericSegmentOffloadMaxSegments_65536K.target
[Unit]
Description=Foo
[Link]
GenericSegmentOffloadMaxSegments=65536K

