# tests/auto/InvalidValue/2.48/bad/Link_TransmitQueues_4097.target
[Unit]
Description=Foo
[Link]
TransmitQueues=4097

