# tests/auto/InvalidValue/2.48/bad/Link_GenericSegmentOffloadMaxBytes_0M.target
[Unit]
Description=Foo
[Link]
GenericSegmentOffloadMaxBytes=0M

