# tests/auto/InvalidValue/2.48/bad/Link_GenericSegmentOffloadMaxSegments_0G.target
[Unit]
Description=Foo
[Link]
GenericSegmentOffloadMaxSegments=0G

