# tests/auto/InvalidValue/2.48/bad/Link_GenericSegmentOffloadMaxSegments_65536G.target
[Unit]
Description=Foo
[Link]
GenericSegmentOffloadMaxSegments=65536G

