# tests/auto/InvalidValue/2.48/bad/Link_ReceiveQueues_0.target
[Unit]
Description=Foo
[Link]
ReceiveQueues=0

