# tests/auto/InvalidValue/2.48/bad/Manager_ManagerEnvironment_1.target
[Unit]
Description=Foo
[Manager]
ManagerEnvironment=1

