# tests/auto/InvalidValue/2.48/bad/Time_ConnectionRetrySec_10000000000000.target
[Unit]
Description=Foo
[Time]
ConnectionRetrySec=10000000000000

