# systemdlint/tests/auto/InvalidValue/2.34/bad/Mount_LogLevelMax_True.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
LogLevelMax=True

