# systemdlint/tests/auto/InvalidValue/2.34/bad/Service_LogLevelMax_True.target
[Unit]
Description=Foo
[Service]
LogLevelMax=True

