# systemdlint/tests/auto/InvalidValue/2.34/bad/Service_LogLevelMax_1.target
[Unit]
Description=Foo
[Service]
LogLevelMax=1

