# systemdlint/tests/auto/InvalidValue/2.34/bad/Socket_LogLevelMax_True.target
[Unit]
Description=Foo
[Socket]
LogLevelMax=True

