# systemdlint/tests/auto/InvalidValue/2.34/bad/Socket_LogLevelMax_1.target
[Unit]
Description=Foo
[Socket]
LogLevelMax=1

