# systemdlint/tests/auto/InvalidValue/2.34/bad/Swap_LogLevelMax_1.target
[Unit]
Description=Foo
[Swap]
LogLevelMax=1

