# systemdlint/tests/auto/InvalidValue/2.34/bad/Swap_LogLevelMax_True.target
[Unit]
Description=Foo
[Swap]
LogLevelMax=True

