# systemdlint/tests/auto/InvalidValue/2.34/bad/Unit_CollectMode_1.target
[Unit]
Description=Foo
Description=abc
CollectMode=1

