# systemdlint/tests/auto/InvalidValue/2.34/bad/Unit_SuccessAction_1.target
[Unit]
Description=Foo
Description=abc
SuccessAction=1

