# systemdlint/tests/auto/InvalidValue/2.34/bad/Unit_CollectMode_True.target
[Unit]
Description=Foo
Description=abc
CollectMode=True

