# systemdlint/tests/auto/InvalidValue/2.40/bad/Mount_MemoryMin_10000000000000M.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
MemoryMin=10000000000000M

