# systemdlint/tests/auto/InvalidValue/2.40/bad/Service_LogRateLimitIntervalSec_True.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=True

