# systemdlint/tests/auto/InvalidValue/2.40/bad/Service_MemoryMin_-1B.target
[Unit]
Description=Foo
[Service]
MemoryMin=-1B

