# systemdlint/tests/auto/InvalidValue/2.40/bad/Service_MemoryMin_-1M.target
[Unit]
Description=Foo
[Service]
MemoryMin=-1M

