# systemdlint/tests/auto/InvalidValue/2.40/bad/Socket_MemoryMin_-1K.target
[Unit]
Description=Foo
[Socket]
MemoryMin=-1K

