# systemdlint/tests/auto/InvalidValue/2.40/bad/Socket_LogRateLimitIntervalSec_1p.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=1p

