# systemdlint/tests/auto/InvalidValue/2.40/bad/Socket_LogRateLimitBurst_3decades.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=3decades

