# systemdlint/tests/auto/InvalidValue/2.40/bad/Socket_LogRateLimitIntervalSec_3decades.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=3decades

