# systemdlint/tests/auto/InvalidValue/2.40/bad/Socket_LogRateLimitIntervalSec_True.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=True

