# systemdlint/tests/auto/InvalidValue/2.40/bad/Socket_MemoryMin_-1B.target
[Unit]
Description=Foo
[Socket]
MemoryMin=-1B

