# systemdlint/tests/auto/InvalidValue/2.40/bad/Socket_LogRateLimitBurst_True.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=True

