# systemdlint/tests/auto/InvalidValue/2.40/bad/Swap_LogRateLimitIntervalSec_True.target
[Unit]
Description=Foo
[Swap]
LogRateLimitIntervalSec=True

