# systemdlint/tests/auto/InvalidValue/2.40/bad/Swap_LogRateLimitBurst_0.5d.target
[Unit]
Description=Foo
[Swap]
LogRateLimitBurst=0.5d

