# systemdlint/tests/auto/InvalidValue/2.40/bad/Swap_LogRateLimitIntervalSec_3decades.target
[Unit]
Description=Foo
[Swap]
LogRateLimitIntervalSec=3decades

