# tests/auto/InvalidValue/2.39/good/Manager_NoNewPrivileges_0.target
[Unit]
Description=Foo
[Manager]
NoNewPrivileges=0

