# tests/auto/InvalidValue/2.39/good/Service_NoNewPrivileges_off.target
[Unit]
Description=Foo
[Service]
NoNewPrivileges=off

