# tests/auto/InvalidValue/2.39/good/Service_NoNewPrivileges_yes.target
[Unit]
Description=Foo
[Service]
NoNewPrivileges=yes

