# tests/auto/InvalidValue/2.39/good/Service_PrivateMounts_1.target
[Unit]
Description=Foo
[Service]
PrivateMounts=1

