# tests/auto/InvalidValue/2.39/good/Service_NoNewPrivileges_no.target
[Unit]
Description=Foo
[Service]
NoNewPrivileges=no

