# tests/auto/InvalidValue/2.39/good/Service_PrivateMounts_yes.target
[Unit]
Description=Foo
[Service]
PrivateMounts=yes

