# tests/auto/InvalidValue/2.39/good/Socket_PrivateMounts_true.target
[Unit]
Description=Foo
[Socket]
PrivateMounts=true

