# tests/auto/InvalidValue/2.39/good/Socket_PrivateMounts_0.target
[Unit]
Description=Foo
[Socket]
PrivateMounts=0

