# tests/auto/InvalidValue/2.39/good/Socket_NoNewPrivileges_false.target
[Unit]
Description=Foo
[Socket]
NoNewPrivileges=false

