# tests/auto/InvalidValue/2.39/good/Socket_PrivateMounts_no.target
[Unit]
Description=Foo
[Socket]
PrivateMounts=no

