# tests/auto/InvalidValue/2.39/good/Swap_PrivateMounts_no.target
[Unit]
Description=Foo
[Swap]
PrivateMounts=no

