# tests/auto/InvalidValue/2.39/good/Swap_NoNewPrivileges_false.target
[Unit]
Description=Foo
[Swap]
NoNewPrivileges=false

