# tests/auto/InvalidValue/2.39/good/Swap_NoNewPrivileges_off.target
[Unit]
Description=Foo
[Swap]
NoNewPrivileges=off

