# tests/auto/InvalidValue/2.39/good/Swap_NoNewPrivileges_true.target
[Unit]
Description=Foo
[Swap]
NoNewPrivileges=true

