# tests/auto/InvalidValue/2.39/good/Swap_NoNewPrivileges_yes.target
[Unit]
Description=Foo
[Swap]
NoNewPrivileges=yes

