# tests/auto/InvalidValue/2.39/good/Swap_PrivateMounts_false.target
[Unit]
Description=Foo
[Swap]
PrivateMounts=false

