# tests/auto/InvalidValue/2.39/good/Swap_NoNewPrivileges_0.target
[Unit]
Description=Foo
[Swap]
NoNewPrivileges=0

