# systemdlint/tests/auto/OptionTooNew/2.27/bad/Mount_PassEnvironment_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
PassEnvironment=123

