# systemdlint/tests/auto/OptionTooNew/2.27/bad/Socket_PassEnvironment_abc.target
[Unit]
Description=Foo
[Socket]
PassEnvironment=123

