# systemdlint/tests/auto/OptionTooNew/2.49/bad/CAKE_SplitGSO_yes.target
[Unit]
Description=Foo
[CAKE]
SplitGSO=123

