# systemdlint/tests/auto/OptionTooNew/2.49/bad/CAKE_Wash_yes.target
[Unit]
Description=Foo
[CAKE]
Wash=123

