# systemdlint/tests/auto/OptionTooNew/2.49/bad/CAN_ClassicDataLengthCode_abc.target
[Unit]
Description=Foo
[CAN]
ClassicDataLengthCode=123

