# systemdlint/tests/auto/OptionTooNew/2.49/bad/CAN_DataPhaseBufferSegment1_0.target
[Unit]
Description=Foo
[CAN]
DataPhaseBufferSegment1=123

