# systemdlint/tests/auto/OptionTooNew/2.49/bad/CAN_PropagationSegment_0.target
[Unit]
Description=Foo
[CAN]
PropagationSegment=123

