# systemdlint/tests/auto/OptionTooNew/2.49/bad/DHCPPrefixDelegation_Assign_yes.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Assign=123

