# systemdlint/tests/auto/OptionTooNew/2.49/bad/IPv6AcceptRA_UseGateway_yes.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseGateway=123

