# systemdlint/tests/auto/OptionTooNew/2.49/bad/Login_HandlePowerKeyLongPress_ignore.target
[Unit]
Description=Foo
[Login]
HandlePowerKeyLongPress=123

