# systemdlint/tests/auto/OptionTooNew/2.49/bad/Manager_DefaultOOMScoreAdjust_-1000.target
[Unit]
Description=Foo
[Manager]
DefaultOOMScoreAdjust=123

