# systemdlint/tests/auto/OptionTooNew/2.49/bad/Mount_ExecSearchPath_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
ExecSearchPath=123

