# systemdlint/tests/auto/OptionTooNew/2.49/bad/Service_AllowedCPUs_abc.target
[Unit]
Description=Foo
[Service]
AllowedCPUs=123

