# systemdlint/tests/auto/OptionTooNew/2.49/bad/Slice_AllowedMemoryNodes_abc.target
[Unit]
Description=Foo
[Slice]
AllowedMemoryNodes=123

