# systemdlint/tests/auto/OptionTooNew/2.49/bad/Socket_RestrictFileSystems_abc.target
[Unit]
Description=Foo
[Socket]
RestrictFileSystems=123

