# systemdlint/tests/auto/OptionTooNew/2.49/bad/Socket_StartupAllowedCPUs_abc.target
[Unit]
Description=Foo
[Socket]
StartupAllowedCPUs=123

