# systemdlint/tests/auto/OptionTooNew/2.49/bad/Swap_TTYRows_0.target
[Unit]
Description=Foo
[Swap]
TTYRows=123

