# systemdlint/tests/auto/OptionTooNew/2.49/bad/WireGuard_AllowedIPs_abc.target
[Unit]
Description=Foo
[WireGuard]
AllowedIPs=123

