# tests/auto/OptionTooNew/2.28/bad/Timer_RandomizedDelaySec_0days.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=123

