# systemdlint/tests/auto/OptionTooNew/2.46/bad/IPv6SendRA_EmitDomains_yes.target
[Unit]
Description=Foo
[IPv6SendRA]
EmitDomains=123

