# systemdlint/tests/auto/OptionTooNew/2.46/bad/OOM_SwapUsedLimit_abc.target
[Unit]
Description=Foo
[OOM]
SwapUsedLimit=123

