# systemdlint/tests/auto/OptionTooNew/2.46/bad/Service_ProcSubset_all.target
[Unit]
Description=Foo
[Service]
ProcSubset=123

