# systemdlint/tests/auto/OptionTooNew/2.46/bad/Socket_FlushPending_yes.target
[Unit]
Description=Foo
[Socket]
FlushPending=123

