# systemdlint/tests/auto/OptionTooNew/2.46/bad/Socket_ProtectProc_noaccess.target
[Unit]
Description=Foo
[Socket]
ProtectProc=123

