# systemdlint/tests/auto/OptionTooNew/2.46/bad/Socket_ProcSubset_all.target
[Unit]
Description=Foo
[Socket]
ProcSubset=123

