# systemdlint/tests/auto/OptionTooNew/2.46/bad/Swap_ManagedOOMMemoryPressure_abc.target
[Unit]
Description=Foo
[Swap]
ManagedOOMMemoryPressure=123

