# systemdlint/tests/auto/OptionTooNew/2.46/bad/Swap_ProcSubset_all.target
[Unit]
Description=Foo
[Swap]
ProcSubset=123

