# systemdlint/tests/auto/OptionTooNew/2.46/bad/Timer_FixedRandomDelay_yes.target
[Unit]
Description=Foo
[Timer]
FixedRandomDelay=123

