# tests/auto/OptionTooNew/2.43/good/DHCPv4_RouteMetric_0.target
[Unit]
Description=Foo
[DHCPv4]
RouteMetric=123

