# tests/auto/OptionTooNew/2.43/good/DHCPv6_RequestBroadcast_abc.target
[Unit]
Description=Foo
[DHCPv6]
RequestBroadcast=123

