# tests/auto/OptionTooNew/2.43/good/DHCPv6_VendorClassIdentifier_abc.target
[Unit]
Description=Foo
[DHCPv6]
VendorClassIdentifier=123

