# tests/auto/OptionTooNew/2.43/good/DHCPv6_DUIDType_abc.target
[Unit]
Description=Foo
[DHCPv6]
DUIDType=123

