# tests/auto/OptionTooNew/2.43/good/DHCPv4_UserClass_abc.target
[Unit]
Description=Foo
[DHCPv4]
UserClass=123

