# tests/auto/OptionTooNew/2.43/good/DHCPv4_ListenPort_1.target
[Unit]
Description=Foo
[DHCPv4]
ListenPort=123

