# tests/auto/OptionTooNew/2.43/good/Socket_NUMAPolicy_default.target
[Unit]
Description=Foo
[Socket]
NUMAPolicy=123

