# tests/auto/OptionTooNew/2.43/good/Socket_DefaultMemoryMin_0K.target
[Unit]
Description=Foo
[Socket]
DefaultMemoryMin=123

