# tests/auto/OptionTooNew/2.43/good/Unit_IPEgressFilterPath_abc.target
[Unit]
Description=Foo
Description=abc
IPEgressFilterPath=123

