# systemdlint/tests/auto/OptionTooNew/2.31/bad/Service_CPUWeight_1.target
[Unit]
Description=Foo
[Service]
CPUWeight=123

