# systemdlint/tests/auto/OptionTooNew/2.31/bad/Service_ProtectControlGroups_yes.target
[Unit]
Description=Foo
[Service]
ProtectControlGroups=123

