# systemdlint/tests/auto/OptionTooNew/2.31/bad/Socket_ProtectControlGroups_yes.target
[Unit]
Description=Foo
[Socket]
ProtectControlGroups=123

