# systemdlint/tests/auto/OptionTooNew/2.31/bad/Socket_MaxConnectionsPerSource_0.target
[Unit]
Description=Foo
[Socket]
MaxConnectionsPerSource=123

