# systemdlint/tests/auto/OptionTooNew/2.31/bad/Swap_DynamicUser_yes.target
[Unit]
Description=Foo
[Swap]
DynamicUser=123

