# systemdlint/tests/auto/OptionTooNew/2.31/bad/Swap_CPUAccounting_yes.target
[Unit]
Description=Foo
[Swap]
CPUAccounting=123

