# tests/auto/OptionDeprecated/2.45/good/Address_PrefixRoute_yes.target
[Unit]
Description=Foo
[Address]
Address=abc
PrefixRoute=123

