# tests/auto/OptionDeprecated/2.45/good/TrafficControlQueueingDiscipline_Parent_root.target
[Unit]
Description=Foo
[TrafficControlQueueingDiscipline]
Parent=123

