# tests/auto/SettingRestricted/2.43/bad/DHCPv4_UseTimezone_yesAnonymize%3Dyes.target
[Unit]
Description=Foo
[DHCPv4]
UseTimezone=yes
Anonymize=yes
