# tests/auto/SettingRestricted/2.43/bad/DHCPv6_UseMTU_yesMTUBytes%3DNone-Anonymize%3Dtrue.target
[Unit]
Description=Foo
[DHCPv6]
UseMTU=yes
MTUBytes=None
Anonymize=true
