# tests/auto/SettingRestricted/2.43/bad/DHCPv4_VendorClassIdentifier_abcAnonymize%3D1.target
[Unit]
Description=Foo
[DHCPv4]
VendorClassIdentifier=abc
Anonymize=1
