# tests/auto/SettingRestricted/2.43/bad/DHCPv6_UseTimezone_yesAnonymize%3Dyes.target
[Unit]
Description=Foo
[DHCPv6]
UseTimezone=yes
Anonymize=yes
