# tests/auto/OptionTooNew/2.12/good/Timer_WakeSystem_yes.target
[Unit]
Description=Foo
[Timer]
WakeSystem=123

