# tests/auto/OptionTooNew/2.49/bad/CAKE_CompensationMode_none.target
[Unit]
Description=Foo
[CAKE]
CompensationMode=123

