# tests/auto/OptionTooNew/2.49/bad/CAKE_MPUBytes_0.target
[Unit]
Description=Foo
[CAKE]
MPUBytes=123

