# tests/auto/OptionTooNew/2.49/bad/CAN_OneShot_yes.target
[Unit]
Description=Foo
[CAN]
OneShot=123

