# tests/auto/OptionTooNew/2.49/bad/CAN_PhaseBufferSegment1_0.target
[Unit]
Description=Foo
[CAN]
PhaseBufferSegment1=123

