# tests/auto/OptionTooNew/2.49/bad/CAN_DataTimeQuantaNSec_0.target
[Unit]
Description=Foo
[CAN]
DataTimeQuantaNSec=123

