# tests/auto/OptionTooNew/2.49/bad/DHCPPrefixDelegation_Token_abc.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Token=123

