# tests/auto/OptionTooNew/2.49/bad/IPoIB_Mode_datagram.target
[Unit]
Description=Foo
[IPoIB]
Mode=123

