# tests/auto/OptionTooNew/2.49/bad/IPoIB_PartitionKey_0.target
[Unit]
Description=Foo
[IPoIB]
PartitionKey=123

