# tests/auto/OptionTooNew/2.49/bad/Login_HandleSuspendKeyLongPress_ignore.target
[Unit]
Description=Foo
[Login]
HandleSuspendKeyLongPress=123

