# tests/auto/OptionTooNew/2.49/bad/Login_HandleHibernateKeyLongPress_ignore.target
[Unit]
Description=Foo
[Login]
HandleHibernateKeyLongPress=123

