# tests/auto/OptionTooNew/2.49/bad/Service_StartupAllowedMemoryNodes_abc.target
[Unit]
Description=Foo
[Service]
StartupAllowedMemoryNodes=123

