# tests/auto/OptionTooNew/2.49/bad/Service_RestrictFileSystems_abc.target
[Unit]
Description=Foo
[Service]
RestrictFileSystems=123

