# tests/auto/OptionTooNew/2.49/bad/Service_StartupAllowedCPUs_abc.target
[Unit]
Description=Foo
[Service]
StartupAllowedCPUs=123

