# tests/auto/OptionTooNew/2.49/bad/Service_TTYColumns_0.target
[Unit]
Description=Foo
[Service]
TTYColumns=123

