# tests/auto/OptionTooNew/2.49/bad/Service_TTYRows_0.target
[Unit]
Description=Foo
[Service]
TTYRows=123

