# tests/auto/OptionTooNew/2.49/bad/Socket_AllowedMemoryNodes_abc.target
[Unit]
Description=Foo
[Socket]
AllowedMemoryNodes=123

