# tests/auto/OptionTooNew/2.49/bad/Swap_AllowedMemoryNodes_abc.target
[Unit]
Description=Foo
[Swap]
AllowedMemoryNodes=123

