# tests/auto/OptionTooNew/2.49/bad/Swap_StartupAllowedMemoryNodes_abc.target
[Unit]
Description=Foo
[Swap]
StartupAllowedMemoryNodes=123

