# tests/auto/OptionTooNew/2.49/bad/Swap_AllowedCPUs_abc.target
[Unit]
Description=Foo
[Swap]
AllowedCPUs=123

