# tests/auto/OptionTooNew/2.49/bad/Unit_AssertIOPressure_abc.target
[Unit]
Description=Foo
Description=abc
AssertIOPressure=123

