# tests/auto/OptionTooNew/2.30/good/Socket_TriggerLimitIntervalSec_0days.target
[Unit]
Description=Foo
[Socket]
TriggerLimitIntervalSec=123

