# tests/auto/MandatoryOptionMissing/2.31/good/L2TPSession_Mandatory_L2TPSession.target
[Unit]
Description=Foo
[L2TPSession]
Name=abc
