# tests/auto/MandatoryOptionMissing/2.31/good/NetDev_Mandatory_NetDev.target
[Unit]
Description=Foo
[NetDev]
Name=abc
