# systemdlint/tests/auto/OptionTooNew/2.34/good/Socket_StandardInputText_abc.target
[Unit]
Description=Foo
[Socket]
StandardInputText=123

