# tests/auto/OptionTooNew/2.27/bad/Swap_PassEnvironment_abc.target
[Unit]
Description=Foo
[Swap]
PassEnvironment=123

