# tests/auto/OptionTooNew/2.50/good/CAN_Loopback_yes.target
[Unit]
Description=Foo
[CAN]
Loopback=123

