# tests/auto/OptionTooNew/2.50/good/CAN_PresumeAck_abc.target
[Unit]
Description=Foo
[CAN]
PresumeAck=123

