# tests/auto/OptionTooNew/2.50/good/CAN_DataPropagationSegment_0.target
[Unit]
Description=Foo
[CAN]
DataPropagationSegment=123

