# tests/auto/OptionTooNew/2.50/good/DHCPPrefixDelegation_UplinkInterface_abc.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
UplinkInterface=123

