# tests/auto/OptionTooNew/2.50/good/DHCPPrefixDelegation_RouteMetric_0.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
RouteMetric=123

