# tests/auto/OptionTooNew/2.50/good/DHCPv6_UseDelegatedPrefix_yes.target
[Unit]
Description=Foo
[DHCPv6]
UseDelegatedPrefix=123

