# tests/auto/OptionTooNew/2.50/good/Mount_RestrictNetworkInterfaces_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
RestrictNetworkInterfaces=123

