# tests/auto/OptionTooNew/2.50/good/Slice_AllowedCPUs_abc.target
[Unit]
Description=Foo
[Slice]
AllowedCPUs=123

