# tests/auto/OptionTooNew/2.50/good/Socket_AllowedCPUs_abc.target
[Unit]
Description=Foo
[Socket]
AllowedCPUs=123

