# tests/auto/OptionTooNew/2.50/good/Socket_ExecSearchPath_abc.target
[Unit]
Description=Foo
[Socket]
ExecSearchPath=123

