# tests/auto/OptionTooNew/2.50/good/Swap_ExecSearchPath_abc.target
[Unit]
Description=Foo
[Swap]
ExecSearchPath=123

