# tests/auto/OptionTooNew/2.50/good/Swap_RestrictFileSystems_abc.target
[Unit]
Description=Foo
[Swap]
RestrictFileSystems=123

