# tests/auto/OptionTooNew/2.50/good/Swap_RestrictNetworkInterfaces_abc.target
[Unit]
Description=Foo
[Swap]
RestrictNetworkInterfaces=123

