# tests/auto/OptionTooNew/2.50/good/Swap_StartupAllowedCPUs_abc.target
[Unit]
Description=Foo
[Swap]
StartupAllowedCPUs=123

