# tests/auto/OptionTooNew/2.50/good/Unit_ConditionCPUPressure_abc.target
[Unit]
Description=Foo
Description=abc
ConditionCPUPressure=123

