# tests/auto/OptionTooNew/2.50/good/Unit_AssertMemoryPressure_abc.target
[Unit]
Description=Foo
Description=abc
AssertMemoryPressure=123

