# systemdlint/tests/auto/SettingRestricted/2.43/good/DHCPv4_SendHostname_yesAnonymize%3Dyes.target
[Unit]
Description=Foo
[DHCPv4]
SendHostname=yes

