# tests/auto/SettingRequires/2.31/bad/Service_GuessMainPID_yesType%3Dforking.target
[Unit]
Description=Foo
[Service]
GuessMainPID=yes

