# tests/auto/SettingRequires/2.31/bad/Socket_SELinuxContext_abcAccept%3Dyes.target
[Unit]
Description=Foo
[Socket]
SELinuxContext=abc

