# tests/auto/SettingRequires/2.31/bad/Tunnel_FOUDestinationPort_1FooOverUDP%3Dyes.target
[Unit]
Description=Foo
[Tunnel]
FOUDestinationPort=1

