# systemdlint/tests/auto/OptionTooNew/2.26/bad/Socket_TasksAccounting_yes.target
[Unit]
Description=Foo
[Socket]
TasksAccounting=123

