# systemdlint/tests/auto/InvalidNumericBase/2.31/bad/CAN_Bitrate_-1K.target
[Unit]
Description=Foo
[CAN]
Bitrate=-1K

