# systemdlint/tests/auto/InvalidNumericBase/2.31/bad/CAN_Bitrate_2001.target
[Unit]
Description=Foo
[CAN]
Bitrate=2001

