# systemdlint/tests/auto/InvalidNumericBase/2.31/bad/CAN_Bitrate_6001.target
[Unit]
Description=Foo
[CAN]
Bitrate=6001

