# systemdlint/tests/auto/InvalidNumericBase/2.31/bad/CAN_Bitrate_3001.target
[Unit]
Description=Foo
[CAN]
Bitrate=3001

