# systemdlint/tests/auto/InvalidNumericBase/2.31/bad/Link_MTUBytes_10000000000000M.target
[Unit]
Description=Foo
[Link]
MACAddress=abc
MTUBytes=10000000000000M

