# systemdlint/tests/auto/InvalidNumericBase/2.31/bad/Link_MTUBytes_2049.target
[Unit]
Description=Foo
[Link]
MACAddress=abc
MTUBytes=2049

