# systemdlint/tests/auto/InvalidNumericBase/2.31/bad/Login_RuntimeDirectorySize_7169.target
[Unit]
Description=Foo
[Login]
RuntimeDirectorySize=7169

