# systemdlint/tests/auto/InvalidNumericBase/2.31/bad/Login_RuntimeDirectorySize_8193.target
[Unit]
Description=Foo
[Login]
RuntimeDirectorySize=8193

