# systemdlint/tests/auto/InvalidNumericBase/2.31/bad/NetDev_MTUBytes_3073.target
[Unit]
Description=Foo
[NetDev]
Name=abc
MTUBytes=3073

