# systemdlint/tests/manual/Multiplicity/2.31/bad/Service_Simple_MultiExecStart.service
[Unit]
Description=multi ExecStart

[Service]
Type=simple
ExecStart=/bin/busybox
ExecStart=/bin/busybox -v