# tests/auto/OptionTooNew/2.42/good/Socket_NetworkNamespacePath_--tmp.target
[Unit]
Description=Foo
[Socket]
NetworkNamespacePath=123

