# tests/auto/OptionTooNew/2.42/good/Swap_NetworkNamespacePath_--tmp.target
[Unit]
Description=Foo
[Swap]
NetworkNamespacePath=123

