# tests/auto/OptionTooNew/2.42/good/Swap_ProtectHostname_yes.target
[Unit]
Description=Foo
[Swap]
ProtectHostname=123

