# tests/auto/OptionTooNew/2.42/good/Unit_ConditionMemory_abc.target
[Unit]
Description=Foo
Description=abc
ConditionMemory=123

