# systemdlint/tests/manual/ReferencedUnitNotFound/2.31/bad/Install_RequiredBy_foo.service.target
[Unit]
Description=Foo

[Install]
RequiredBy=foo.service