# tests/auto/OptionTooNew/2.34/bad/Mount_UnsetEnvironment_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
UnsetEnvironment=123

