# tests/auto/OptionTooNew/2.34/bad/Service_StateDirectoryMode_0.target
[Unit]
Description=Foo
[Service]
StateDirectoryMode=123

