# tests/auto/OptionTooNew/2.34/bad/Service_IPAddressDeny_1.1.1.1--0.target
[Unit]
Description=Foo
[Service]
IPAddressDeny=123

