# tests/auto/OptionTooNew/2.34/bad/Socket_LockPersonality_yes.target
[Unit]
Description=Foo
[Socket]
LockPersonality=123

