# tests/auto/OptionTooNew/2.34/bad/Socket_UnsetEnvironment_abc.target
[Unit]
Description=Foo
[Socket]
UnsetEnvironment=123

