# tests/auto/OptionTooNew/2.34/bad/Socket_KeyringMode_inherit.target
[Unit]
Description=Foo
[Socket]
KeyringMode=123

