# tests/auto/OptionTooNew/2.34/bad/Socket_StateDirectoryMode_0.target
[Unit]
Description=Foo
[Socket]
StateDirectoryMode=123

