# tests/auto/OptionTooNew/2.34/bad/Swap_CacheDirectory_abc.target
[Unit]
Description=Foo
[Swap]
CacheDirectory=123

