# tests/auto/OptionTooNew/2.34/bad/Swap_RuntimeDirectoryMode_0.target
[Unit]
Description=Foo
[Swap]
RuntimeDirectoryMode=123

