# tests/auto/InvalidValue/2.44/good/DHCPServer_EmitSIP_false.target
[Unit]
Description=Foo
[DHCPServer]
EmitSIP=false

