# tests/auto/InvalidValue/2.44/good/DHCPServer_SendOption_abc.target
[Unit]
Description=Foo
[DHCPServer]
SendOption=abc

