# tests/auto/InvalidValue/2.44/good/DHCPServer_EmitSIP_true.target
[Unit]
Description=Foo
[DHCPServer]
EmitSIP=true

