# tests/auto/InvalidValue/2.44/good/DHCPServer_EmitSIP_yes.target
[Unit]
Description=Foo
[DHCPServer]
EmitSIP=yes

