# tests/auto/InvalidValue/2.44/good/DHCPServer_EmitSIP_0.target
[Unit]
Description=Foo
[DHCPServer]
EmitSIP=0

