# tests/auto/InvalidValue/2.44/good/IPv6RoutePrefix_Route_abc.target
[Unit]
Description=Foo
[IPv6RoutePrefix]
Route=abc

