# tests/auto/InvalidValue/2.44/good/LifetimeSec_Route_9999999999999.target
[Unit]
Description=Foo
[LifetimeSec]
Route=9999999999999

