# tests/auto/InvalidValue/2.44/good/LifetimeSec_Route_0.target
[Unit]
Description=Foo
[LifetimeSec]
Route=0

