# tests/auto/InvalidValue/2.44/good/Match_WLANInterfaceType_ad-hoc.target
[Unit]
Description=Foo
[Match]
WLANInterfaceType=ad-hoc

