# tests/auto/InvalidValue/2.44/good/Match_WLANInterfaceType_wds.target
[Unit]
Description=Foo
[Match]
WLANInterfaceType=wds

