# tests/auto/InvalidValue/2.44/good/Match_WLANInterfaceType_nan.target
[Unit]
Description=Foo
[Match]
WLANInterfaceType=nan

