# tests/auto/InvalidValue/2.44/good/Mount_RestartKillSignal_34.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
RestartKillSignal=34

