# tests/auto/InvalidValue/2.44/good/Route_IPServiceType_CS6.target
[Unit]
Description=Foo
[Route]
IPServiceType=CS6

