# tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_7minutes.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=7minutes

