# tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_25years.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=25years

