# tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_11m.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=11m

