# tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_20usec.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=20usec

