# tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_10M.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=10M

