# tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_21us.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=21us

