# tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_27y.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=27y

