# tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_6h.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=6h

