# tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_0days.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=0days

