# tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_2d.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=2d

