# tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_16seconds.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=16seconds

