# tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_23week.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=23week

