# tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_12month.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=12month

