# tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_1day.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=1day

