# tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_8minute.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=8minute

