# tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_15ms.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=15ms

