# tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_3hours.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=3hours

