# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGFPE.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGFPE

