# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGINT.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGINT

