# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGIOT.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGIOT

