# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_32.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=32

