# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGPROF.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGPROF

