# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_38.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=38

