# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGCHLD.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGCHLD

