# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_11.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=11

