# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_DATAERR.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=DATAERR

