# tests/auto/InvalidValue/2.44/good/Service_ProtectKernelLogs_off.target
[Unit]
Description=Foo
[Service]
ProtectKernelLogs=off

