# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGUSR2.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGUSR2

