# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_15.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=15

