# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGURG.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGURG

