# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_8.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=8

