# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGPIPE.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGPIPE

