# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGBUS.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGBUS

