# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGEMT.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGEMT

