# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGSTKFLT.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGSTKFLT

