# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGSEGV.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGSEGV

