# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGTTOU.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGTTOU

