# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGTERM.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGTERM

