# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGWINCH.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGWINCH

