# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGUNUSED.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGUNUSED

