# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGLOST.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGLOST

