# tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGHUP.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGHUP

