# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGPROF.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGPROF

