# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGXCPU.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGXCPU

