# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGCLD.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGCLD

