# tests/auto/InvalidValue/2.44/good/Socket_ProtectKernelLogs_false.target
[Unit]
Description=Foo
[Socket]
ProtectKernelLogs=false

