# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGVTALRM.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGVTALRM

