# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGCHLD.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGCHLD

