# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGHUP.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGHUP

