# tests/auto/InvalidValue/2.44/good/Socket_ProtectKernelLogs_true.target
[Unit]
Description=Foo
[Socket]
ProtectKernelLogs=true

