# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_25.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=25

