# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGTERM.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGTERM

