# tests/auto/InvalidValue/2.44/good/Socket_ProtectKernelLogs_0.target
[Unit]
Description=Foo
[Socket]
ProtectKernelLogs=0

