# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGFPE.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGFPE

