# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGALRM.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGALRM

