# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_22.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=22

