# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGPOLL.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGPOLL

