# tests/auto/InvalidValue/2.44/good/Socket_ProtectKernelLogs_1.target
[Unit]
Description=Foo
[Socket]
ProtectKernelLogs=1

