# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGINT.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGINT

