# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGIOT.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGIOT

