# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGPIPE.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGPIPE

