# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGTTOU.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGTTOU

