# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGEMT.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGEMT

