# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGQUIT.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGQUIT

