# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_16.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=16

