# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGXFSZ.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGXFSZ

