# tests/auto/InvalidValue/2.44/good/Socket_ProtectKernelLogs_yes.target
[Unit]
Description=Foo
[Socket]
ProtectKernelLogs=yes

