# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_23.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=23

