# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGABRT.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGABRT

