# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGSYS.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGSYS

