# tests/auto/InvalidValue/2.44/good/Socket_ProtectKernelLogs_off.target
[Unit]
Description=Foo
[Socket]
ProtectKernelLogs=off

