# tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGKILL.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGKILL

