# tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGTERM.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGTERM

