# tests/auto/InvalidValue/2.44/good/Swap_ProtectKernelLogs_on.target
[Unit]
Description=Foo
[Swap]
ProtectKernelLogs=on

