# tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGKILL.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGKILL

