# tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGSTKFLT.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGSTKFLT

