# tests/auto/InvalidValue/2.44/good/Swap_ProtectKernelLogs_yes.target
[Unit]
Description=Foo
[Swap]
ProtectKernelLogs=yes

