# tests/auto/InvalidValue/2.44/good/Swap_ProtectKernelLogs_no.target
[Unit]
Description=Foo
[Swap]
ProtectKernelLogs=no

