# tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGALRM.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGALRM

