# tests/auto/InvalidValue/2.44/good/Swap_ProtectKernelLogs_0.target
[Unit]
Description=Foo
[Swap]
ProtectKernelLogs=0

