# tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGEMT.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGEMT

