# tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_52.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=52

