# tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGABRT.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGABRT

