# tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGPROF.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGPROF

