# tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_9.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=9

