# tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGHUP.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGHUP

