# tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGSEGV.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGSEGV

