# tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_36.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=36

