# tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGPWR.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGPWR

