# tests/auto/InvalidValue/2.44/good/Swap_ProtectKernelLogs_off.target
[Unit]
Description=Foo
[Swap]
ProtectKernelLogs=off

