# tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGCONT.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGCONT

