# tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGCHLD.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGCHLD

