# tests/auto/InvalidValue/2.44/good/TrafficControlQueueingDiscipline_NetworkEmulatorDelayJitterSec_9999999999999.target
[Unit]
Description=Foo
[TrafficControlQueueingDiscipline]
NetworkEmulatorDelayJitterSec=9999999999999

