# tests/auto/InvalidValue/2.44/good/TrafficControlQueueingDiscipline_NetworkEmulatorDelaySec_9999999999999.target
[Unit]
Description=Foo
[TrafficControlQueueingDiscipline]
NetworkEmulatorDelaySec=9999999999999

