# tests/auto/InvalidValue/2.44/good/TrafficControlQueueingDiscipline_Parent_root.target
[Unit]
Description=Foo
[TrafficControlQueueingDiscipline]
Parent=root

