# systemdlint/tests/auto/InvalidValue/2.50/bad/CAKE_AutoRateIngress_foo.target
[Unit]
Description=Foo
[CAKE]
AutoRateIngress=foo

