# systemdlint/tests/auto/InvalidValue/2.50/bad/CAKE_AutoRateIngress_bar.target
[Unit]
Description=Foo
[CAKE]
AutoRateIngress=bar

