# systemdlint/tests/auto/InvalidValue/2.50/bad/CAKE_Wash_foo.target
[Unit]
Description=Foo
[CAKE]
Wash=foo

