# systemdlint/tests/auto/InvalidValue/2.50/bad/CAKE_SplitGSO_foo.target
[Unit]
Description=Foo
[CAKE]
SplitGSO=foo

