# systemdlint/tests/auto/InvalidValue/2.50/bad/CAKE_Wash_bar.target
[Unit]
Description=Foo
[CAKE]
Wash=bar

