# systemdlint/tests/auto/InvalidValue/2.50/bad/CAKE_UseRawPacketSize_bar.target
[Unit]
Description=Foo
[CAKE]
UseRawPacketSize=bar

