# systemdlint/tests/auto/InvalidValue/2.50/bad/CAKE_SplitGSO_bar.target
[Unit]
Description=Foo
[CAKE]
SplitGSO=bar

