# systemdlint/tests/auto/InvalidValue/2.50/bad/CAN_DataTimeQuantaNSec_-1.target
[Unit]
Description=Foo
[CAN]
DataTimeQuantaNSec=-1

