# systemdlint/tests/auto/InvalidValue/2.50/bad/CAN_PhaseBufferSegment2_4294967296.target
[Unit]
Description=Foo
[CAN]
PhaseBufferSegment2=4294967296

