# systemdlint/tests/auto/InvalidValue/2.50/bad/CAN_SyncJumpWidth_4294967296.target
[Unit]
Description=Foo
[CAN]
SyncJumpWidth=4294967296

