# systemdlint/tests/auto/InvalidValue/2.50/bad/CAN_PropagationSegment_-1.target
[Unit]
Description=Foo
[CAN]
PropagationSegment=-1

