# systemdlint/tests/auto/InvalidValue/2.50/bad/CAN_DataPhaseBufferSegment2_4294967296.target
[Unit]
Description=Foo
[CAN]
DataPhaseBufferSegment2=4294967296

