# systemdlint/tests/auto/InvalidValue/2.50/bad/CAN_Loopback_baz.target
[Unit]
Description=Foo
[CAN]
Loopback=baz

