# systemdlint/tests/auto/InvalidValue/2.50/bad/DHCPPrefixDelegation_Announce_bar.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Announce=bar

