# systemdlint/tests/auto/InvalidValue/2.50/bad/DHCPPrefixDelegation_Assign_baz.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Assign=baz

