# systemdlint/tests/auto/InvalidValue/2.50/bad/DHCPPrefixDelegation_Announce_foo.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Announce=foo

