# systemdlint/tests/auto/InvalidValue/2.50/bad/DHCPPrefixDelegation_Announce_baz.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Announce=baz

