# systemdlint/tests/auto/InvalidValue/2.50/bad/DHCPPrefixDelegation_Assign_foo.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Assign=foo

