# systemdlint/tests/auto/InvalidValue/2.50/bad/DHCPPrefixDelegation_Assign_bar.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Assign=bar

