# systemdlint/tests/auto/InvalidValue/2.50/bad/IPv6AcceptRA_UseMTU_baz.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseMTU=baz

