# systemdlint/tests/auto/InvalidValue/2.50/bad/Mount_TTYColumns_-1.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
TTYColumns=-1

