# systemdlint/tests/auto/InvalidValue/2.50/bad/Scope_RuntimeRandomizedExtraSec_3decades.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=3decades

