# systemdlint/tests/auto/InvalidValue/2.50/bad/Service_RuntimeRandomizedExtraSec_0.5d.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=0.5d

