# systemdlint/tests/auto/InvalidValue/2.50/bad/Service_RuntimeRandomizedExtraSec_1p.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=1p

