# systemdlint/tests/auto/InvalidValue/2.50/bad/Service_TTYColumns_-1.target
[Unit]
Description=Foo
[Service]
TTYColumns=-1

