# systemdlint/tests/auto/InvalidValue/2.50/bad/Swap_TTYColumns_-1.target
[Unit]
Description=Foo
[Swap]
TTYColumns=-1

