# systemdlint/tests/auto/InvalidValue/2.50/bad/Swap_TTYColumns_10000000000000.target
[Unit]
Description=Foo
[Swap]
TTYColumns=10000000000000

