# systemdlint/tests/auto/InvalidValue/2.50/bad/Time_SaveIntervalSec_True.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=True

