# systemdlint/tests/auto/SettingRestricted/2.43/bad/DHCPv4_UseMTU_yesAnonymize%3Dtrue-MTUBytes%3DNone.target
[Unit]
Description=Foo
[DHCPv4]
UseMTU=yes
Anonymize=true
MTUBytes=None
