# systemdlint/tests/auto/SettingRestricted/2.43/bad/DHCPv4_UseMTU_yesMTUBytes%3DNone-Anonymize%3D1.target
[Unit]
Description=Foo
[DHCPv4]
UseMTU=yes
MTUBytes=None
Anonymize=1
