# systemdlint/tests/auto/SettingRestricted/2.43/bad/DHCPv6_SendHostname_yesAnonymize%3Dyes.target
[Unit]
Description=Foo
[DHCPv6]
SendHostname=yes
Anonymize=yes
