# systemdlint/tests/auto/SettingRestricted/2.43/bad/DHCPv6_VendorClassIdentifier_abcAnonymize%3D1.target
[Unit]
Description=Foo
[DHCPv6]
VendorClassIdentifier=abc
Anonymize=1
