# systemdlint/tests/auto/SettingRestricted/2.43/bad/DHCPv6_ClientIdentifier_macAnonymize%3D1.target
[Unit]
Description=Foo
[DHCPv6]
ClientIdentifier=mac
Anonymize=1
