# systemdlint/tests/auto/OptionDeprecated/2.46/good/IPv6PrefixDelegation_Domains_abc.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
Domains=123

