# systemdlint/tests/auto/OptionTooNew/2.35/good/Service_IPAccounting_yes.target
[Unit]
Description=Foo
[Service]
IPAccounting=123

