# systemdlint/tests/auto/OptionTooNew/2.35/good/Service_ConfigurationDirectoryMode_0.target
[Unit]
Description=Foo
[Service]
ConfigurationDirectoryMode=123

