# systemdlint/tests/auto/OptionTooNew/2.35/good/Socket_IPAddressAllow_1.1.1.1--0.target
[Unit]
Description=Foo
[Socket]
IPAddressAllow=123

