# systemdlint/tests/auto/OptionTooNew/2.35/good/Socket_RuntimeDirectoryPreserve_yes.target
[Unit]
Description=Foo
[Socket]
RuntimeDirectoryPreserve=123

