# systemdlint/tests/auto/OptionTooNew/2.35/good/Swap_ConfigurationDirectory_abc.target
[Unit]
Description=Foo
[Swap]
ConfigurationDirectory=123

