# tests/auto/SettingRestricted/2.31/bad/DHCP_VendorClassIdentifier_abcAnonymize%3D1.target
[Unit]
Description=Foo
[DHCP]
VendorClassIdentifier=abc
Anonymize=1
