# tests/auto/SettingRestricted/2.31/bad/DHCP_UseMTU_yesMTUBytes%3DNone-Anonymize%3D1.target
[Unit]
Description=Foo
[DHCP]
UseMTU=yes
MTUBytes=None
Anonymize=1
