# tests/auto/SettingRestricted/2.31/bad/DHCP_UseTimezone_yesAnonymize%3Dyes.target
[Unit]
Description=Foo
[DHCP]
UseTimezone=yes
Anonymize=yes
