# tests/auto/SettingRestricted/2.31/bad/DHCP_SendHostname_yesAnonymize%3D1.target
[Unit]
Description=Foo
[DHCP]
SendHostname=yes
Anonymize=1
