# tests/auto/SettingRestricted/2.31/bad/Network_ActiveSlave_yesMode%3Dbalance-alb.target
[Unit]
Description=Foo
[Network]
ActiveSlave=yes
Mode=balance-alb
