# systemdlint/tests/auto/SettingRequires/2.31/good/BridgeVLAN_PVID_abcVLAN%3DNone.target
[Unit]
Description=Foo
[BridgeVLAN]
PVID=abc
VLAN=None
