# systemdlint/tests/auto/SettingRequires/2.31/good/BridgeVLAN_EgressUntagged_abcVLAN%3DNone.target
[Unit]
Description=Foo
[BridgeVLAN]
EgressUntagged=abc
VLAN=None
