# systemdlint/tests/auto/SettingRequires/2.31/good/DHCP_Anonymize_yesMACAddressPolicy%3D1.target
[Unit]
Description=Foo
[DHCP]
Anonymize=yes
MACAddressPolicy=1
