# systemdlint/tests/auto/SettingRequires/2.31/good/DHCP_Anonymize_yesMACAddressPolicy%3Dtrue.target
[Unit]
Description=Foo
[DHCP]
Anonymize=yes
MACAddressPolicy=true
