# systemdlint/tests/auto/SettingRequires/2.31/good/Link_AllMulticast_yesMulticast%3Dtrue.target
[Unit]
Description=Foo
[Link]
AllMulticast=yes
Multicast=true
