# systemdlint/tests/auto/SettingRequires/2.31/good/Login_IdleActionSec_0MIdleAction%3DNone.target
[Unit]
Description=Foo
[Login]
IdleActionSec=0M
IdleAction=None
