# systemdlint/tests/auto/SettingRequires/2.31/good/Network_IPMasquerade_yesIPForward%3Dipv4.target
[Unit]
Description=Foo
[Network]
IPMasquerade=yes
IPForward=ipv4
