# systemdlint/tests/auto/SettingRequires/2.31/good/Service_BusName_abcType%3Ddbus.target
[Unit]
Description=Foo
[Service]
BusName=abc
Type=dbus
