# systemdlint/tests/auto/SettingRequires/2.31/good/Service_GuessMainPID_yesType%3Dforking.target
[Unit]
Description=Foo
[Service]
GuessMainPID=yes
Type=forking
