# systemdlint/tests/auto/SettingRequires/2.31/good/Service_RootDirectoryStartOnly_yesExecStart%3DNone.target
[Unit]
Description=Foo
[Service]
RootDirectoryStartOnly=yes
ExecStart=None
