# systemdlint/tests/auto/SettingRequires/2.31/good/Socket_ListenUSBFunction_abcUSBFunctionStrings%3DNone-USBFunctionDescriptors%3DNone.target
[Unit]
Description=Foo
[Socket]
ListenUSBFunction=abc
USBFunctionStrings=None
USBFunctionDescriptors=None
