# systemdlint/tests/auto/SettingRequires/2.31/good/Socket_Service_syslog.serviceAccept%3Dno.target
[Unit]
Description=Foo
[Socket]
Service=syslog.service
Accept=no
