# systemdlint/tests/auto/OptionTooNew/2.45/good/ControlledDelay_IntervalSec_0days.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=123

