# systemdlint/tests/auto/OptionTooNew/2.45/good/FairQueueing_InitialQuantum_0K.target
[Unit]
Description=Foo
[FairQueueing]
InitialQuantum=123

