# systemdlint/tests/auto/OptionTooNew/2.45/good/FairQueueingControlledDelay_TargetSec_0days.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=123

