# systemdlint/tests/auto/OptionTooNew/2.45/good/FairQueueing_Handle_00.target
[Unit]
Description=Foo
[FairQueueing]
Handle=123

