# systemdlint/tests/auto/OptionTooNew/2.45/good/FairQueueing_PacketLimit_0.target
[Unit]
Description=Foo
[FairQueueing]
PacketLimit=123

