# systemdlint/tests/auto/OptionTooNew/2.45/good/FairQueueing_Quantum_0K.target
[Unit]
Description=Foo
[FairQueueing]
Quantum=123

