# systemdlint/tests/auto/OptionTooNew/2.45/good/FairQueueingControlledDelay_CEThresholdSec_0days.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=123

