# systemdlint/tests/auto/OptionTooNew/2.45/good/FairQueueingControlledDelay_ECN_yes.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
ECN=123

