# systemdlint/tests/auto/OptionTooNew/2.45/good/FairQueueingControlledDelay_Quantum_0K.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
Quantum=123

