# systemdlint/tests/auto/OptionTooNew/2.45/good/Link_AlternativeName_abc.target
[Unit]
Description=Foo
[Link]
AlternativeName=123

