# systemdlint/tests/auto/OptionTooNew/2.45/good/Match_PermanentMACAddress_abc.target
[Unit]
Description=Foo
[Match]
PermanentMACAddress=123

