# systemdlint/tests/auto/OptionTooNew/2.45/good/QDics_Handle_00.target
[Unit]
Description=Foo
[QDics]
Handle=123

