# systemdlint/tests/auto/OptionTooNew/2.45/good/Socket_ProtectClock_yes.target
[Unit]
Description=Foo
[Socket]
ProtectClock=123

