# systemdlint/tests/auto/OptionTooNew/2.45/good/TokenBucketFilter_PeakRate_0K.target
[Unit]
Description=Foo
[TokenBucketFilter]
PeakRate=123

