# systemdlint/tests/auto/OptionTooNew/2.45/good/TokenBucketFilter_Rate_0K.target
[Unit]
Description=Foo
[TokenBucketFilter]
Rate=123

