# systemdlint/tests/auto/OptionTooNew/2.45/good/TokenBucketFilter_MTUBytes_0K.target
[Unit]
Description=Foo
[TokenBucketFilter]
MTUBytes=123

