# systemdlint/tests/auto/OptionTooNew/2.45/good/TrivialLinkEqualizer_Parent_root.target
[Unit]
Description=Foo
[TrivialLinkEqualizer]
Parent=123

