# tests/auto/InvalidValue/2.42/bad/Mount_NetworkNamespacePath_1.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
NetworkNamespacePath=1

