# tests/auto/InvalidValue/2.42/bad/Service_NetworkNamespacePath_True.target
[Unit]
Description=Foo
[Service]
NetworkNamespacePath=True

