# tests/auto/InvalidValue/2.42/bad/Service_CPUQuotaPeriodSec_1p.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=1p

