# tests/auto/InvalidValue/2.42/bad/Service_CPUQuotaPeriodSec_2potatoes.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=2potatoes

