# tests/auto/InvalidValue/2.42/bad/Service_CPUQuotaPeriodSec_0.5d.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=0.5d

