# tests/auto/InvalidValue/2.42/bad/Service_ProtectHostname_baz.target
[Unit]
Description=Foo
[Service]
ProtectHostname=baz

