# tests/auto/InvalidValue/2.42/bad/Socket_ProtectHostname_foo.target
[Unit]
Description=Foo
[Socket]
ProtectHostname=foo

