# tests/auto/InvalidValue/2.42/bad/Socket_NetworkNamespacePath_1.target
[Unit]
Description=Foo
[Socket]
NetworkNamespacePath=1

