# tests/auto/InvalidValue/2.42/bad/Socket_ProtectHostname_bar.target
[Unit]
Description=Foo
[Socket]
ProtectHostname=bar

