# tests/auto/InvalidValue/2.42/bad/Swap_NetworkNamespacePath_True.target
[Unit]
Description=Foo
[Swap]
NetworkNamespacePath=True

