# tests/auto/InvalidValue/2.42/bad/Swap_NetworkNamespacePath_1.target
[Unit]
Description=Foo
[Swap]
NetworkNamespacePath=1

