# tests/auto/InvalidValue/2.42/bad/Swap_CPUQuotaPeriodSec_True.target
[Unit]
Description=Foo
[Swap]
CPUQuotaPeriodSec=True

