# tests/auto/InvalidValue/2.42/bad/Swap_ProtectHostname_baz.target
[Unit]
Description=Foo
[Swap]
ProtectHostname=baz

