# tests/auto/InvalidValue/2.42/bad/Swap_ProtectHostname_foo.target
[Unit]
Description=Foo
[Swap]
ProtectHostname=foo

