# tests/auto/InvalidValue/2.42/bad/Swap_NetworkNamespacePath_tmp.target
[Unit]
Description=Foo
[Swap]
NetworkNamespacePath=tmp

