# tests/auto/InvalidValue/2.42/bad/Swap_CPUQuotaPeriodSec_0.5d.target
[Unit]
Description=Foo
[Swap]
CPUQuotaPeriodSec=0.5d

