# tests/auto/InvalidValue/2.50/good/CAKE_SplitGSO_yes.target
[Unit]
Description=Foo
[CAKE]
SplitGSO=yes

