# tests/auto/InvalidValue/2.50/good/CAKE_NAT_false.target
[Unit]
Description=Foo
[CAKE]
NAT=false

