# tests/auto/InvalidValue/2.50/good/CAKE_PriorityQueueingPreset_diffserv3.target
[Unit]
Description=Foo
[CAKE]
PriorityQueueingPreset=diffserv3

