# tests/auto/InvalidValue/2.50/good/CAKE_AutoRateIngress_0.target
[Unit]
Description=Foo
[CAKE]
AutoRateIngress=0

