# tests/auto/InvalidValue/2.50/good/CAKE_FlowIsolationMode_triple.target
[Unit]
Description=Foo
[CAKE]
FlowIsolationMode=triple

