# tests/auto/InvalidValue/2.50/good/CAKE_AutoRateIngress_true.target
[Unit]
Description=Foo
[CAKE]
AutoRateIngress=true

