# tests/auto/InvalidValue/2.50/good/CAKE_SplitGSO_on.target
[Unit]
Description=Foo
[CAKE]
SplitGSO=on

