# tests/auto/InvalidValue/2.50/good/CAKE_SplitGSO_off.target
[Unit]
Description=Foo
[CAKE]
SplitGSO=off

