# tests/auto/InvalidValue/2.50/good/CAKE_UseRawPacketSize_false.target
[Unit]
Description=Foo
[CAKE]
UseRawPacketSize=false

