# tests/auto/InvalidValue/2.50/good/CAKE_SplitGSO_0.target
[Unit]
Description=Foo
[CAKE]
SplitGSO=0

