# tests/auto/InvalidValue/2.50/good/CAKE_AutoRateIngress_on.target
[Unit]
Description=Foo
[CAKE]
AutoRateIngress=on

