# tests/auto/InvalidValue/2.50/good/CAKE_SplitGSO_no.target
[Unit]
Description=Foo
[CAKE]
SplitGSO=no

