# tests/auto/InvalidValue/2.50/good/CAKE_FlowIsolationMode_dual-src-host.target
[Unit]
Description=Foo
[CAKE]
FlowIsolationMode=dual-src-host

