# tests/auto/InvalidValue/2.50/good/CAKE_PriorityQueueingPreset_precedence.target
[Unit]
Description=Foo
[CAKE]
PriorityQueueingPreset=precedence

