# tests/auto/InvalidValue/2.50/good/CAKE_CompensationMode_ptm.target
[Unit]
Description=Foo
[CAKE]
CompensationMode=ptm

