# tests/auto/InvalidValue/2.50/good/CAKE_UseRawPacketSize_on.target
[Unit]
Description=Foo
[CAKE]
UseRawPacketSize=on

