# tests/auto/InvalidValue/2.50/good/CAKE_Wash_yes.target
[Unit]
Description=Foo
[CAKE]
Wash=yes

