# tests/auto/InvalidValue/2.50/good/CAKE_MPUBytes_256.target
[Unit]
Description=Foo
[CAKE]
MPUBytes=256

