# tests/auto/InvalidValue/2.50/good/CAKE_AutoRateIngress_no.target
[Unit]
Description=Foo
[CAKE]
AutoRateIngress=no

