# tests/auto/InvalidValue/2.50/good/CAKE_UseRawPacketSize_1.target
[Unit]
Description=Foo
[CAKE]
UseRawPacketSize=1

