# tests/auto/InvalidValue/2.50/good/CAKE_AutoRateIngress_off.target
[Unit]
Description=Foo
[CAKE]
AutoRateIngress=off

