# tests/auto/InvalidValue/2.50/good/CAKE_UseRawPacketSize_no.target
[Unit]
Description=Foo
[CAKE]
UseRawPacketSize=no

