# tests/auto/InvalidValue/2.50/good/CAKE_Wash_false.target
[Unit]
Description=Foo
[CAKE]
Wash=false

