# tests/auto/InvalidValue/2.50/good/CAKE_UseRawPacketSize_true.target
[Unit]
Description=Foo
[CAKE]
UseRawPacketSize=true

