# tests/auto/InvalidValue/2.50/good/CAN_PropagationSegment_4294967295.target
[Unit]
Description=Foo
[CAN]
PropagationSegment=4294967295

