# tests/auto/InvalidValue/2.50/good/CAN_Loopback_no.target
[Unit]
Description=Foo
[CAN]
Loopback=no

