# tests/auto/InvalidValue/2.50/good/CAN_TimeQuantaNSec_4294967295.target
[Unit]
Description=Foo
[CAN]
TimeQuantaNSec=4294967295

