# tests/auto/InvalidValue/2.50/good/CAN_PhaseBufferSegment2_0.target
[Unit]
Description=Foo
[CAN]
PhaseBufferSegment2=0

