# tests/auto/InvalidValue/2.50/good/CAN_Loopback_true.target
[Unit]
Description=Foo
[CAN]
Loopback=true

