# tests/auto/InvalidValue/2.50/good/CAN_DataPropagationSegment_4294967295.target
[Unit]
Description=Foo
[CAN]
DataPropagationSegment=4294967295

