# tests/auto/InvalidValue/2.50/good/CAN_Loopback_on.target
[Unit]
Description=Foo
[CAN]
Loopback=on

