# tests/auto/InvalidValue/2.50/good/CAN_DataPhaseBufferSegment1_0.target
[Unit]
Description=Foo
[CAN]
DataPhaseBufferSegment1=0

