# tests/auto/InvalidValue/2.50/good/CAN_OneShot_false.target
[Unit]
Description=Foo
[CAN]
OneShot=false

