# tests/auto/InvalidValue/2.50/good/CAN_Loopback_0.target
[Unit]
Description=Foo
[CAN]
Loopback=0

