# tests/auto/InvalidValue/2.50/good/CAN_OneShot_on.target
[Unit]
Description=Foo
[CAN]
OneShot=on

