# tests/auto/InvalidValue/2.50/good/CAN_DataSyncJumpWidth_0.target
[Unit]
Description=Foo
[CAN]
DataSyncJumpWidth=0

