# tests/auto/InvalidValue/2.50/good/DHCPPrefixDelegation_Assign_no.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Assign=no

