# tests/auto/InvalidValue/2.50/good/DHCPPrefixDelegation_Announce_true.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Announce=true

