# tests/auto/InvalidValue/2.50/good/DHCPPrefixDelegation_Assign_on.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Assign=on

