# tests/auto/InvalidValue/2.50/good/DHCPPrefixDelegation_Token_abc.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Token=abc

