# tests/auto/InvalidValue/2.50/good/DHCPPrefixDelegation_Assign_1.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Assign=1

