# tests/auto/InvalidValue/2.50/good/DHCPPrefixDelegation_Announce_no.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Announce=no

