# tests/auto/InvalidValue/2.50/good/DHCPv6_UseDelegatedPrefix_false.target
[Unit]
Description=Foo
[DHCPv6]
UseDelegatedPrefix=false

