# tests/auto/InvalidValue/2.50/good/DHCPPrefixDelegation_Announce_on.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Announce=on

