# tests/auto/InvalidValue/2.50/good/DHCPPrefixDelegation_ManageTemporaryAddress_abc.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
ManageTemporaryAddress=abc

