# tests/auto/InvalidValue/2.50/good/DHCPPrefixDelegation_Announce_0.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Announce=0

