# tests/auto/InvalidValue/2.50/good/IPoIB_Mode_connected.target
[Unit]
Description=Foo
[IPoIB]
Mode=connected

