# tests/auto/InvalidValue/2.50/good/IPoIB_IgnoreUserspaceMulticastGroup_1.target
[Unit]
Description=Foo
[IPoIB]
IgnoreUserspaceMulticastGroup=1

