# tests/auto/InvalidValue/2.50/good/IPv6AcceptRA_UseMTU_no.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseMTU=no

