# tests/auto/InvalidValue/2.50/good/IPv6AcceptRA_UseGateway_yes.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseGateway=yes

