# tests/auto/InvalidValue/2.50/good/IPv6AcceptRA_UseGateway_false.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseGateway=false

