# tests/auto/InvalidValue/2.50/good/Login_HandleHibernateKeyLongPress_lock.target
[Unit]
Description=Foo
[Login]
HandleHibernateKeyLongPress=lock

