# tests/auto/InvalidValue/2.50/good/Login_HandleHibernateKeyLongPress_ignore.target
[Unit]
Description=Foo
[Login]
HandleHibernateKeyLongPress=ignore

