# tests/auto/InvalidValue/2.50/good/Login_HandlePowerKeyLongPress_kexec.target
[Unit]
Description=Foo
[Login]
HandlePowerKeyLongPress=kexec

