# tests/auto/InvalidValue/2.50/good/Login_HandleRebootKeyLongPress_suspend-then-hibernate.target
[Unit]
Description=Foo
[Login]
HandleRebootKeyLongPress=suspend-then-hibernate

