# tests/auto/InvalidValue/2.50/good/Login_HandleSuspendKeyLongPress_poweroff.target
[Unit]
Description=Foo
[Login]
HandleSuspendKeyLongPress=poweroff

