# tests/auto/InvalidValue/2.50/good/Login_HandleHibernateKeyLongPress_poweroff.target
[Unit]
Description=Foo
[Login]
HandleHibernateKeyLongPress=poweroff

