# tests/auto/InvalidValue/2.50/good/Login_HandleSuspendKeyLongPress_suspend-then-hibernate.target
[Unit]
Description=Foo
[Login]
HandleSuspendKeyLongPress=suspend-then-hibernate

