# tests/auto/InvalidValue/2.50/good/Login_HandleRebootKeyLongPress_kexec.target
[Unit]
Description=Foo
[Login]
HandleRebootKeyLongPress=kexec

