# tests/auto/InvalidValue/2.50/good/Login_HandlePowerKeyLongPress_poweroff.target
[Unit]
Description=Foo
[Login]
HandlePowerKeyLongPress=poweroff

