# tests/auto/InvalidValue/2.50/good/Login_HandlePowerKeyLongPress_suspend-then-hibernate.target
[Unit]
Description=Foo
[Login]
HandlePowerKeyLongPress=suspend-then-hibernate

