# tests/auto/InvalidValue/2.50/good/Login_HandleHibernateKeyLongPress_hibernate.target
[Unit]
Description=Foo
[Login]
HandleHibernateKeyLongPress=hibernate

