# tests/auto/InvalidValue/2.50/good/Login_HandleRebootKeyLongPress_poweroff.target
[Unit]
Description=Foo
[Login]
HandleRebootKeyLongPress=poweroff

