# tests/auto/InvalidValue/2.50/good/Login_HandlePowerKeyLongPress_reboot.target
[Unit]
Description=Foo
[Login]
HandlePowerKeyLongPress=reboot

