# tests/auto/InvalidValue/2.50/good/Login_HandleSuspendKeyLongPress_kexec.target
[Unit]
Description=Foo
[Login]
HandleSuspendKeyLongPress=kexec

