# tests/auto/InvalidValue/2.50/good/Login_HandlePowerKeyLongPress_hybrid-sleep.target
[Unit]
Description=Foo
[Login]
HandlePowerKeyLongPress=hybrid-sleep

