# tests/auto/InvalidValue/2.50/good/Login_HandleHibernateKeyLongPress_kexec.target
[Unit]
Description=Foo
[Login]
HandleHibernateKeyLongPress=kexec

