# tests/auto/InvalidValue/2.50/good/Mount_StartupAllowedMemoryNodes_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
StartupAllowedMemoryNodes=abc

