# tests/auto/InvalidValue/2.50/good/Path_TriggerLimitIntervalSec_9min.target
[Unit]
Description=Foo
[Path]
TriggerLimitIntervalSec=9min

