# tests/auto/InvalidValue/2.50/good/Path_TriggerLimitIntervalSec_16seconds.target
[Unit]
Description=Foo
[Path]
TriggerLimitIntervalSec=16seconds

