# tests/auto/InvalidValue/2.50/good/Path_TriggerLimitIntervalSec_12month.target
[Unit]
Description=Foo
[Path]
TriggerLimitIntervalSec=12month

