# tests/auto/InvalidValue/2.50/good/Path_TriggerLimitBurst_9999999999999.target
[Unit]
Description=Foo
[Path]
TriggerLimitBurst=9999999999999

