# tests/auto/InvalidValue/2.50/good/Path_TriggerLimitIntervalSec_26year.target
[Unit]
Description=Foo
[Path]
TriggerLimitIntervalSec=26year

