# tests/auto/InvalidValue/2.50/good/Path_TriggerLimitIntervalSec_infinity.target
[Unit]
Description=Foo
[Path]
TriggerLimitIntervalSec=infinity

