# tests/auto/InvalidValue/2.50/good/Path_TriggerLimitIntervalSec_6h.target
[Unit]
Description=Foo
[Path]
TriggerLimitIntervalSec=6h

