# tests/auto/InvalidValue/2.50/good/Path_TriggerLimitIntervalSec_28%C2%B5s.target
[Unit]
Description=Foo
[Path]
TriggerLimitIntervalSec=28µs

