# tests/auto/InvalidValue/2.50/good/Path_TriggerLimitIntervalSec_22weeks.target
[Unit]
Description=Foo
[Path]
TriggerLimitIntervalSec=22weeks

