# tests/auto/InvalidValue/2.50/good/RoutingPolicyRule_SuppressInterfaceGroup_2147483647.target
[Unit]
Description=Foo
[RoutingPolicyRule]
SuppressInterfaceGroup=2147483647

