# tests/auto/InvalidValue/2.50/good/RoutingPolicyRule_SuppressInterfaceGroup_0.target
[Unit]
Description=Foo
[RoutingPolicyRule]
SuppressInterfaceGroup=0

