# tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_9min.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=9min

