# tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_12month.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=12month

