# tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_13months.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=13months

