# tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_22weeks.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=22weeks

