# tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_23week.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=23week

