# tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_7minutes.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=7minutes

