# tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_15ms.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=15ms

