# tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_14msec.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=14msec

