# tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_1day.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=1day

