# tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_27y.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=27y

