# tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_11m.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=11m

