# tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_8minute.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=8minute

