# tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_19s.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=19s

