# tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_24w.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=24w

