# tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_6h.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=6h

