# tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_21us.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=21us

