# tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_22weeks.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=22weeks

