# tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_19s.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=19s

