# tests/auto/InvalidValue/2.50/good/Service_ExitType_main.target
[Unit]
Description=Foo
[Service]
ExitType=main

