# tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_4hour.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=4hour

