# tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_infinity.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=infinity

