# tests/auto/InvalidValue/2.50/good/Service_AllowedCPUs_abc.target
[Unit]
Description=Foo
[Service]
AllowedCPUs=abc

