# tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_16seconds.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=16seconds

