# tests/auto/InvalidValue/2.50/good/Service_RestrictFileSystems_abc.target
[Unit]
Description=Foo
[Service]
RestrictFileSystems=abc

