# tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_8minute.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=8minute

