# tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_11m.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=11m

