# tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_28%C2%B5s.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=28µs

