# tests/auto/InvalidValue/2.50/good/Service_TTYColumns_0.target
[Unit]
Description=Foo
[Service]
TTYColumns=0

