# tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_17second.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=17second

