# tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_12month.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=12month

