# tests/auto/InvalidValue/2.50/good/Service_ExecSearchPath_abc.target
[Unit]
Description=Foo
[Service]
ExecSearchPath=abc

