# tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_2d.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=2d

