# tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_25years.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=25years

