# tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_14msec.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=14msec

