# tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_10M.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=10M

