# tests/auto/InvalidValue/2.50/good/Service_TTYColumns_9999999999999.target
[Unit]
Description=Foo
[Service]
TTYColumns=9999999999999

