# tests/auto/InvalidValue/2.50/good/Service_TTYRows_0.target
[Unit]
Description=Foo
[Service]
TTYRows=0

