# tests/auto/InvalidValue/2.50/good/Slice_RestrictNetworkInterfaces_abc.target
[Unit]
Description=Foo
[Slice]
RestrictNetworkInterfaces=abc

