# tests/auto/InvalidValue/2.50/good/Slice_StartupAllowedMemoryNodes_abc.target
[Unit]
Description=Foo
[Slice]
StartupAllowedMemoryNodes=abc

