# tests/auto/InvalidValue/2.50/good/Slice_StartupAllowedCPUs_abc.target
[Unit]
Description=Foo
[Slice]
StartupAllowedCPUs=abc

