# tests/auto/InvalidValue/2.50/good/Slice_AllowedMemoryNodes_abc.target
[Unit]
Description=Foo
[Slice]
AllowedMemoryNodes=abc

