# tests/auto/InvalidValue/2.50/good/Socket_RestrictNetworkInterfaces_abc.target
[Unit]
Description=Foo
[Socket]
RestrictNetworkInterfaces=abc

