# tests/auto/InvalidValue/2.50/good/Socket_StartupAllowedCPUs_abc.target
[Unit]
Description=Foo
[Socket]
StartupAllowedCPUs=abc

