# tests/auto/InvalidValue/2.50/good/Socket_StartupAllowedMemoryNodes_abc.target
[Unit]
Description=Foo
[Socket]
StartupAllowedMemoryNodes=abc

