# tests/auto/InvalidValue/2.50/good/Socket_RestrictFileSystems_abc.target
[Unit]
Description=Foo
[Socket]
RestrictFileSystems=abc

