# tests/auto/InvalidValue/2.50/good/Swap_AllowedMemoryNodes_abc.target
[Unit]
Description=Foo
[Swap]
AllowedMemoryNodes=abc

