# tests/auto/InvalidValue/2.50/good/Swap_RestrictFileSystems_abc.target
[Unit]
Description=Foo
[Swap]
RestrictFileSystems=abc

