# tests/auto/InvalidValue/2.50/good/Swap_StartupAllowedMemoryNodes_abc.target
[Unit]
Description=Foo
[Swap]
StartupAllowedMemoryNodes=abc

