# tests/auto/InvalidValue/2.50/good/Swap_AllowedCPUs_abc.target
[Unit]
Description=Foo
[Swap]
AllowedCPUs=abc

