# tests/auto/InvalidValue/2.50/good/Swap_TTYColumns_0.target
[Unit]
Description=Foo
[Swap]
TTYColumns=0

