# tests/auto/InvalidValue/2.50/good/Swap_StartupAllowedCPUs_abc.target
[Unit]
Description=Foo
[Swap]
StartupAllowedCPUs=abc

