# tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_20usec.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=20usec

