# tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_6h.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=6h

