# tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_14msec.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=14msec

