# tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_27y.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=27y

