# tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_23week.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=23week

