# tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_15ms.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=15ms

