# tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_1day.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=1day

