# tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_25years.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=25years

