# tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_10M.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=10M

