# tests/auto/InvalidValue/2.50/good/Unit_ConditionIOPressure_abc.target
[Unit]
Description=Foo
Description=abc
ConditionIOPressure=abc

