# tests/auto/InvalidValue/2.50/good/Unit_AssertIOPressure_abc.target
[Unit]
Description=Foo
Description=abc
AssertIOPressure=abc

