# tests/auto/InvalidValue/2.50/good/Unit_AssertCPUPressure_abc.target
[Unit]
Description=Foo
Description=abc
AssertCPUPressure=abc

