# tests/auto/InvalidValue/2.50/good/WireGuardPeer_AllowedIPs_abc.target
[Unit]
Description=Foo
[WireGuardPeer]
AllowedIPs=abc

