# tests/auto/InvalidValue/2.50/good/WireGuard_AllowedIPs_abc.target
[Unit]
Description=Foo
[WireGuard]
AllowedIPs=abc

