# systemdlint/tests/auto/OptionTooNew/2.10/bad/Manager_DefaultMemoryAccounting_yes.target
[Unit]
Description=Foo
[Manager]
DefaultMemoryAccounting=123

