# tests/auto/InvalidValue/2.34/good/Service_LogLevelMax_debug.target
[Unit]
Description=Foo
[Service]
LogLevelMax=debug

