# tests/auto/InvalidValue/2.34/good/Service_StandardInputText_abc.target
[Unit]
Description=Foo
[Service]
StandardInputText=abc

