# tests/auto/InvalidValue/2.34/good/Service_StandardInput_null.target
[Unit]
Description=Foo
[Service]
StandardInput=null

