# tests/auto/InvalidValue/2.34/good/Service_StandardInputData_abc.target
[Unit]
Description=Foo
[Service]
StandardInputData=abc

