# tests/auto/InvalidValue/2.34/good/Service_LogLevelMax_alert.target
[Unit]
Description=Foo
[Service]
LogLevelMax=alert

