# tests/auto/InvalidValue/2.34/good/Service_StandardInput_data.target
[Unit]
Description=Foo
[Service]
StandardInput=data

