# tests/auto/InvalidValue/2.34/good/Socket_StandardInput_data.target
[Unit]
Description=Foo
[Socket]
StandardInput=data

