# tests/auto/InvalidValue/2.34/good/Socket_StandardInputText_abc.target
[Unit]
Description=Foo
[Socket]
StandardInputText=abc

