# tests/auto/InvalidValue/2.34/good/Socket_LogLevelMax_alert.target
[Unit]
Description=Foo
[Socket]
LogLevelMax=alert

