# tests/auto/InvalidValue/2.34/good/Socket_StandardInput_null.target
[Unit]
Description=Foo
[Socket]
StandardInput=null

