# tests/auto/InvalidValue/2.34/good/Socket_LogLevelMax_info.target
[Unit]
Description=Foo
[Socket]
LogLevelMax=info

