# tests/auto/InvalidValue/2.34/good/Socket_LogLevelMax_emerg.target
[Unit]
Description=Foo
[Socket]
LogLevelMax=emerg

