# tests/auto/InvalidValue/2.34/good/Socket_LogLevelMax_notice.target
[Unit]
Description=Foo
[Socket]
LogLevelMax=notice

