# tests/auto/InvalidValue/2.34/good/Swap_StandardInput_socket.target
[Unit]
Description=Foo
[Swap]
StandardInput=socket

