# tests/auto/InvalidValue/2.34/good/Swap_StandardInputData_abc.target
[Unit]
Description=Foo
[Swap]
StandardInputData=abc

