# tests/auto/InvalidValue/2.34/good/Swap_LogLevelMax_alert.target
[Unit]
Description=Foo
[Swap]
LogLevelMax=alert

