# tests/auto/InvalidValue/2.34/good/Swap_StandardInput_abc.target
[Unit]
Description=Foo
[Swap]
StandardInput=abc

