# tests/auto/InvalidValue/2.34/good/Swap_StandardInput_tty-fail.target
[Unit]
Description=Foo
[Swap]
StandardInput=tty-fail

