# tests/auto/InvalidValue/2.34/good/Swap_StandardInput_data.target
[Unit]
Description=Foo
[Swap]
StandardInput=data

