# tests/auto/InvalidValue/2.34/good/Swap_LogLevelMax_notice.target
[Unit]
Description=Foo
[Swap]
LogLevelMax=notice

