# tests/auto/InvalidValue/2.34/good/Swap_LogLevelMax_info.target
[Unit]
Description=Foo
[Swap]
LogLevelMax=info

