# tests/auto/InvalidValue/2.34/good/Unit_SuccessAction_poweroff-force.target
[Unit]
Description=Foo
Description=abc
SuccessAction=poweroff-force

