# tests/auto/OptionDeprecated/2.29/bad/Unit_RequiresOverridable_abc.target
[Unit]
Description=Foo
Description=abc
RequiresOverridable=123

