# systemdlint/tests/auto/OptionTooNew/2.39/good/Swap_PrivateMounts_yes.target
[Unit]
Description=Foo
[Swap]
PrivateMounts=123

