# tests/auto/OptionTooNew/2.40/good/Socket_MemoryMin_0K.target
[Unit]
Description=Foo
[Socket]
MemoryMin=123

