# tests/auto/OptionTooNew/2.40/good/Unit_SuccessActionExitStatus_0.target
[Unit]
Description=Foo
Description=abc
SuccessActionExitStatus=123

