# systemdlint/tests/manual/ConflictingOptions/2.31/good/Unit_RequiredBy.target
[Unit]
Description=Before and After
RequiredBy=foo.service