# systemdlint/tests/manual/ConflictingOptions/2.31/good/Unit_WantedBy.target
[Unit]
Description=Before and After
WantedBy=foo.service