# systemdlint/tests/manual/ConflictingOptions/2.31/good/Unit_Before.target
[Unit]
Description=Before and After
Before=foo.service