# systemdlint/tests/auto/OptionTooNew/2.17/bad/Service_SmackProcessLabel_abc.target
[Unit]
Description=Foo
[Service]
SmackProcessLabel=123

