# systemdlint/tests/auto/InvalidValue/2.45/bad/Address_AddPrefixRoute_baz.target
[Unit]
Description=Foo
[Address]
Address=abc
AddPrefixRoute=baz

