# systemdlint/tests/auto/InvalidValue/2.45/bad/ControlledDelay_IntervalSec_2potatoes.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=2potatoes

