# systemdlint/tests/auto/InvalidValue/2.45/bad/ControlledDelay_CEThresholdSec_0.5d.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=0.5d

