# systemdlint/tests/auto/InvalidValue/2.45/bad/ControlledDelay_TargetSec_0.5d.target
[Unit]
Description=Foo
[ControlledDelay]
TargetSec=0.5d

