# systemdlint/tests/auto/InvalidValue/2.45/bad/ControlledDelay_CEThresholdSec_True.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=True

