# systemdlint/tests/auto/InvalidValue/2.45/bad/ControlledDelay_IntervalSec_1p.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=1p

