# systemdlint/tests/auto/InvalidValue/2.45/bad/ControlledDelay_TargetSec_True.target
[Unit]
Description=Foo
[ControlledDelay]
TargetSec=True

