# systemdlint/tests/auto/InvalidValue/2.45/bad/ControlledDelay_IntervalSec_True.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=True

