# systemdlint/tests/auto/InvalidValue/2.45/bad/ControlledDelay_IntervalSec_0.5d.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=0.5d

