# systemdlint/tests/auto/InvalidValue/2.45/bad/ControlledDelay_TargetSec_2potatoes.target
[Unit]
Description=Foo
[ControlledDelay]
TargetSec=2potatoes

