# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueing_OrphanMask_-1.target
[Unit]
Description=Foo
[FairQueueing]
OrphanMask=-1

