# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueingControlledDelay_TargetSec_0.5d.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=0.5d

