# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueing_FlowLimit_10000000000000.target
[Unit]
Description=Foo
[FairQueueing]
FlowLimit=10000000000000

