# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueing_Pacing_baz.target
[Unit]
Description=Foo
[FairQueueing]
Pacing=baz

