# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueing_MaximumRate_-1M.target
[Unit]
Description=Foo
[FairQueueing]
MaximumRate=-1M

