# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueingControlledDelay_Quantum_-1K.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
Quantum=-1K

