# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueing_MaximumRate_-1G.target
[Unit]
Description=Foo
[FairQueueing]
MaximumRate=-1G

