# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueingControlledDelay_Quantum_-1M.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
Quantum=-1M

