# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueingControlledDelay_Quantum_10000000000000G.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
Quantum=10000000000000G

