# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueing_FlowLimit_-1.target
[Unit]
Description=Foo
[FairQueueing]
FlowLimit=-1

