# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueing_InitialQuantum_-1K.target
[Unit]
Description=Foo
[FairQueueing]
InitialQuantum=-1K

