# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueing_CEThresholdSec_2potatoes.target
[Unit]
Description=Foo
[FairQueueing]
CEThresholdSec=2potatoes

