# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueing_Buckets_10000000000000.target
[Unit]
Description=Foo
[FairQueueing]
Buckets=10000000000000

