# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueing_MaximumRate_10000000000000M.target
[Unit]
Description=Foo
[FairQueueing]
MaximumRate=10000000000000M

