# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueing_Pacing_foo.target
[Unit]
Description=Foo
[FairQueueing]
Pacing=foo

