# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueing_MaximumRate_-1K.target
[Unit]
Description=Foo
[FairQueueing]
MaximumRate=-1K

