# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueing_InitialQuantum_10000000000000M.target
[Unit]
Description=Foo
[FairQueueing]
InitialQuantum=10000000000000M

