# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueing_PacketLimit_10000000000000.target
[Unit]
Description=Foo
[FairQueueing]
PacketLimit=10000000000000

