# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueingControlledDelay_CEThresholdSec_3decades.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=3decades

