# systemdlint/tests/auto/InvalidValue/2.45/bad/FairQueueing_OrphanMask_10000000000000.target
[Unit]
Description=Foo
[FairQueueing]
OrphanMask=10000000000000

