# systemdlint/tests/auto/InvalidValue/2.45/bad/NetworkEmulator_LossRate_-1%25.target
[Unit]
Description=Foo
[NetworkEmulator]
LossRate=-1%

