# systemdlint/tests/auto/InvalidValue/2.45/bad/NetworkEmulator_DelayJitterSec_0.5d.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=0.5d

