# systemdlint/tests/auto/InvalidValue/2.45/bad/NetworkEmulator_DelayJitterSec_1p.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=1p

