# systemdlint/tests/auto/InvalidValue/2.45/bad/NetworkEmulator_DelayJitterSec_True.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=True

