# systemdlint/tests/auto/InvalidValue/2.45/bad/NetworkEmulator_DelaySec_-1.target
[Unit]
Description=Foo
[NetworkEmulator]
DelaySec=-1

