# systemdlint/tests/auto/InvalidValue/2.45/bad/RoutingPolicyRule_SuppressPrefixLength_129.target
[Unit]
Description=Foo
[RoutingPolicyRule]
SuppressPrefixLength=129

