# systemdlint/tests/auto/InvalidValue/2.45/bad/Service_ProtectClock_bar.target
[Unit]
Description=Foo
[Service]
ProtectClock=bar

