# systemdlint/tests/auto/InvalidValue/2.45/bad/Socket_ProtectClock_foo.target
[Unit]
Description=Foo
[Socket]
ProtectClock=foo

