# systemdlint/tests/auto/InvalidValue/2.45/bad/StochasticFairnessQueueing_PerturbPeriodSec_True.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=True

