# systemdlint/tests/auto/InvalidValue/2.45/bad/StochasticFairnessQueueing_PerturbPeriodSec_0.5d.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=0.5d

