# systemdlint/tests/auto/InvalidValue/2.45/bad/StochasticFairnessQueueing_PerturbPeriodSec_1p.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=1p

