# systemdlint/tests/auto/InvalidValue/2.45/bad/Swap_ProtectClock_foo.target
[Unit]
Description=Foo
[Swap]
ProtectClock=foo

