# systemdlint/tests/auto/InvalidValue/2.45/bad/Swap_ProtectClock_baz.target
[Unit]
Description=Foo
[Swap]
ProtectClock=baz

