# systemdlint/tests/auto/InvalidValue/2.45/bad/TokenBucketFilter_MTUBytes_10000000000000K.target
[Unit]
Description=Foo
[TokenBucketFilter]
MTUBytes=10000000000000K

