# systemdlint/tests/auto/InvalidValue/2.45/bad/TokenBucketFilter_LatencySec_1p.target
[Unit]
Description=Foo
[TokenBucketFilter]
LatencySec=1p

