# systemdlint/tests/auto/InvalidValue/2.45/bad/TokenBucketFilter_Burst_10000000000000M.target
[Unit]
Description=Foo
[TokenBucketFilter]
Burst=10000000000000M

