# systemdlint/tests/auto/InvalidValue/2.45/bad/TokenBucketFilter_MPUBytes_10000000000000K.target
[Unit]
Description=Foo
[TokenBucketFilter]
MPUBytes=10000000000000K

