# systemdlint/tests/auto/InvalidValue/2.45/bad/TokenBucketFilter_PeakRate_-1K.target
[Unit]
Description=Foo
[TokenBucketFilter]
PeakRate=-1K

