# systemdlint/tests/auto/InvalidValue/2.45/bad/TokenBucketFilter_LatencySec_3decades.target
[Unit]
Description=Foo
[TokenBucketFilter]
LatencySec=3decades

