# systemdlint/tests/auto/InvalidValue/2.45/bad/TrivialLinkEqualizer_Id_10000000000000.target
[Unit]
Description=Foo
[TrivialLinkEqualizer]
Id=10000000000000

