# tests/auto/OptionTooNew/2.33/good/Service_RestrictNamespaces_abc.target
[Unit]
Description=Foo
[Service]
RestrictNamespaces=123

