# tests/auto/OptionTooNew/2.33/good/Service_BindReadOnlyPaths_abc.target
[Unit]
Description=Foo
[Service]
BindReadOnlyPaths=123

