# tests/auto/OptionTooNew/2.33/good/Socket_BindPaths_abc.target
[Unit]
Description=Foo
[Socket]
BindPaths=123

