# tests/auto/OptionTooNew/2.33/good/Swap_BindPaths_abc.target
[Unit]
Description=Foo
[Swap]
BindPaths=123

