# tests/auto/OptionTooNew/2.33/good/Swap_RestrictNamespaces_abc.target
[Unit]
Description=Foo
[Swap]
RestrictNamespaces=123

