# tests/auto/InvalidValue/2.38/good/Service_TemporaryFileSystem_abc.target
[Unit]
Description=Foo
[Service]
TemporaryFileSystem=abc

