# tests/auto/OptionTooNew/2.44/bad/ControlledDelay_Handle_00.target
[Unit]
Description=Foo
[ControlledDelay]
Handle=123

