# tests/auto/OptionTooNew/2.44/bad/FairQueueingControlledDelay_Parent_root.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
Parent=123

