# tests/auto/OptionTooNew/2.44/bad/FairQueueingControlledDelay_PacketLimit_0.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
PacketLimit=123

