# tests/auto/OptionTooNew/2.44/bad/NetworkEmulator_LossRate_0%25.target
[Unit]
Description=Foo
[NetworkEmulator]
LossRate=123

