# tests/auto/OptionTooNew/2.44/bad/Service_ProtectClock_yes.target
[Unit]
Description=Foo
[Service]
ProtectClock=123

