# tests/auto/OptionTooNew/2.44/bad/StochasticFairnessQueueing_Parent_root.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
Parent=123

