# tests/auto/OptionTooNew/2.44/bad/TokenBucketFilter_Parent_root.target
[Unit]
Description=Foo
[TokenBucketFilter]
Parent=123

