# tests/auto/OptionTooNew/2.44/bad/TokenBucketFilter_MPUBytes_0K.target
[Unit]
Description=Foo
[TokenBucketFilter]
MPUBytes=123

