# tests/auto/OptionTooNew/2.44/bad/TokenBucketFilter_Handle_00.target
[Unit]
Description=Foo
[TokenBucketFilter]
Handle=123

