# tests/auto/OptionTooNew/2.44/bad/TokenBucketFilter_LatencySec_0days.target
[Unit]
Description=Foo
[TokenBucketFilter]
LatencySec=123

