# systemdlint/tests/auto/OptionTooNew/2.47/bad/IPv6AcceptRA_RouteAllowList_abc.target
[Unit]
Description=Foo
[IPv6AcceptRA]
RouteAllowList=123

