# systemdlint/tests/auto/OptionTooNew/2.47/bad/Link_GenericSegmentOffloadMaxSegments_1K.target
[Unit]
Description=Foo
[Link]
GenericSegmentOffloadMaxSegments=123

