# systemdlint/tests/auto/OptionTooNew/2.47/bad/Link_TransmitQueues_1.target
[Unit]
Description=Foo
[Link]
TransmitQueues=123

