# systemdlint/tests/auto/OptionTooNew/2.47/bad/RoutingPolicyRule_Type_blackhole.target
[Unit]
Description=Foo
[RoutingPolicyRule]
Type=123

