# systemdlint/tests/auto/OptionTooNew/2.47/bad/Service_ExecPaths_abc.target
[Unit]
Description=Foo
[Service]
ExecPaths=123

