# systemdlint/tests/auto/OptionTooNew/2.47/bad/Swap_ExecPaths_abc.target
[Unit]
Description=Foo
[Swap]
ExecPaths=123

