# tests/auto/OptionDeprecated/2.43/bad/DHCP_CriticalConnection_yes.target
[Unit]
Description=Foo
[DHCP]
CriticalConnection=123

