# tests/auto/OptionDeprecated/2.43/bad/Neighbor_MACAddress_01%3A23%3A45%3A67%3A89%3Aab%2000-11-22-33-44-55%20AABB.CCDD.EEFF.target
[Unit]
Description=Foo
[Neighbor]
MACAddress=123

