# tests/auto/InvalidValue/2.46/good/BFIFO_PacketLimit_0.target
[Unit]
Description=Foo
[BFIFO]
PacketLimit=0

