# tests/auto/InvalidValue/2.46/good/BFIFO_Handle_0x1.target
[Unit]
Description=Foo
[BFIFO]
Handle=0x1

