# tests/auto/InvalidValue/2.46/good/BFIFO_Parent_root.target
[Unit]
Description=Foo
[BFIFO]
Parent=root

