# tests/auto/InvalidValue/2.46/good/BFIFO_Parent_ingress.target
[Unit]
Description=Foo
[BFIFO]
Parent=ingress

