# tests/auto/InvalidValue/2.46/good/CAKE_OverheadBytes_-64.target
[Unit]
Description=Foo
[CAKE]
OverheadBytes=-64

