# tests/auto/InvalidValue/2.46/good/CAKE_Parent_abc.target
[Unit]
Description=Foo
[CAKE]
Parent=abc

