# tests/auto/InvalidValue/2.46/good/CAKE_Parent_ingress.target
[Unit]
Description=Foo
[CAKE]
Parent=ingress

