# tests/auto/InvalidValue/2.46/good/CAN_FDNonISO_0.target
[Unit]
Description=Foo
[CAN]
FDNonISO=0

