# tests/auto/InvalidValue/2.46/good/CAN_FDMode_on.target
[Unit]
Description=Foo
[CAN]
FDMode=on

