# tests/auto/InvalidValue/2.46/good/CAN_Termination_yes.target
[Unit]
Description=Foo
[CAN]
Termination=yes

