# tests/auto/InvalidValue/2.46/good/CAN_ListenOnly_false.target
[Unit]
Description=Foo
[CAN]
ListenOnly=false

