# tests/auto/InvalidValue/2.46/good/CAN_BitRate_1M.target
[Unit]
Description=Foo
[CAN]
BitRate=1M

