# tests/auto/InvalidValue/2.46/good/CAN_ListenOnly_true.target
[Unit]
Description=Foo
[CAN]
ListenOnly=true

