# tests/auto/InvalidValue/2.46/good/CAN_FDNonISO_on.target
[Unit]
Description=Foo
[CAN]
FDNonISO=on

