# tests/auto/InvalidValue/2.46/good/CAN_FDMode_no.target
[Unit]
Description=Foo
[CAN]
FDMode=no

