# tests/auto/InvalidValue/2.46/good/CAN_FDNonISO_off.target
[Unit]
Description=Foo
[CAN]
FDNonISO=off

