# tests/auto/InvalidValue/2.46/good/CAN_DataSamplePoint_abc.target
[Unit]
Description=Foo
[CAN]
DataSamplePoint=abc

