# tests/auto/InvalidValue/2.46/good/CAN_FDMode_false.target
[Unit]
Description=Foo
[CAN]
FDMode=false

