# tests/auto/InvalidValue/2.46/good/CAN_ListenOnly_no.target
[Unit]
Description=Foo
[CAN]
ListenOnly=no

