# tests/auto/InvalidValue/2.46/good/CAN_Termination_off.target
[Unit]
Description=Foo
[CAN]
Termination=off

