# tests/auto/InvalidValue/2.46/good/CAN_FDMode_yes.target
[Unit]
Description=Foo
[CAN]
FDMode=yes

