# tests/auto/InvalidValue/2.46/good/CAN_BitRate_1K.target
[Unit]
Description=Foo
[CAN]
BitRate=1K

