# tests/auto/InvalidValue/2.46/good/CAN_BitRate_4294967295M.target
[Unit]
Description=Foo
[CAN]
BitRate=4294967295M

