# tests/auto/InvalidValue/2.46/good/CAN_Termination_on.target
[Unit]
Description=Foo
[CAN]
Termination=on

