# tests/auto/InvalidValue/2.46/good/CAN_ListenOnly_0.target
[Unit]
Description=Foo
[CAN]
ListenOnly=0

