# tests/auto/InvalidValue/2.46/good/DeficitRoundRobinSchedulerClass_QuantumBytes_9999999999999G.target
[Unit]
Description=Foo
[DeficitRoundRobinSchedulerClass]
QuantumBytes=9999999999999G

