# tests/auto/InvalidValue/2.46/good/DeficitRoundRobinSchedulerClass_ClassId_abc.target
[Unit]
Description=Foo
[DeficitRoundRobinSchedulerClass]
ClassId=abc

