# tests/auto/InvalidValue/2.46/good/DeficitRoundRobinSchedulerClass_Parent_abc.target
[Unit]
Description=Foo
[DeficitRoundRobinSchedulerClass]
Parent=abc

