# tests/auto/InvalidValue/2.46/good/DHCPServer_EmitSMTP_true.target
[Unit]
Description=Foo
[DHCPServer]
EmitSMTP=true

