# tests/auto/InvalidValue/2.46/good/DHCPServer_EmitLPR_yes.target
[Unit]
Description=Foo
[DHCPServer]
EmitLPR=yes

