# tests/auto/InvalidValue/2.46/good/DHCPPrefixDelegation_SubnetId_0x0.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
SubnetId=0x0

