# tests/auto/InvalidValue/2.46/good/DHCPv6PrefixDelegation_Assign_no.target
[Unit]
Description=Foo
[DHCPv6PrefixDelegation]
Assign=no

