# tests/auto/InvalidValue/2.46/good/DHCPServer_POP3_abc.target
[Unit]
Description=Foo
[DHCPServer]
POP3=abc

