# tests/auto/InvalidValue/2.46/good/DHCPServer_EmitLPR_0.target
[Unit]
Description=Foo
[DHCPServer]
EmitLPR=0

