# tests/auto/InvalidValue/2.46/good/DHCPServer_SMTP_abc.target
[Unit]
Description=Foo
[DHCPServer]
SMTP=abc

