# tests/auto/InvalidValue/2.46/good/DHCPServer_SendVendorOption_abc.target
[Unit]
Description=Foo
[DHCPServer]
SendVendorOption=abc

