# tests/auto/InvalidValue/2.46/good/DHCPv6_SendOption_abc.target
[Unit]
Description=Foo
[DHCPv6]
SendOption=abc

