# tests/auto/InvalidValue/2.46/good/DHCPv4_FallbackLeaseLifetimeSec_forever.target
[Unit]
Description=Foo
[DHCPv4]
FallbackLeaseLifetimeSec=forever

