# tests/auto/InvalidValue/2.46/good/DHCPServer_EmitSMTP_0.target
[Unit]
Description=Foo
[DHCPServer]
EmitSMTP=0

