# tests/auto/InvalidValue/2.46/good/DHCPServer_EmitPOP3_false.target
[Unit]
Description=Foo
[DHCPServer]
EmitPOP3=false

