# tests/auto/InvalidValue/2.46/good/DHCPv6PrefixDelegation_SubnetId_0x7fffffffffffffff.target
[Unit]
Description=Foo
[DHCPv6PrefixDelegation]
SubnetId=0x7fffffffffffffff

