# tests/auto/InvalidValue/2.46/good/DHCPServer_EmitLPR_false.target
[Unit]
Description=Foo
[DHCPServer]
EmitLPR=false

