# tests/auto/InvalidValue/2.46/good/DHCPv4_UseGateway_no.target
[Unit]
Description=Foo
[DHCPv4]
UseGateway=no

