# tests/auto/InvalidValue/2.46/good/DHCPServer_LRP_abc.target
[Unit]
Description=Foo
[DHCPServer]
LRP=abc

