# tests/auto/InvalidValue/2.46/good/DHCPServer_EmitLPR_on.target
[Unit]
Description=Foo
[DHCPServer]
EmitLPR=on

