# tests/auto/InvalidValue/2.46/good/DHCPv6PrefixDelegation_Assign_false.target
[Unit]
Description=Foo
[DHCPv6PrefixDelegation]
Assign=false

