# tests/auto/InvalidValue/2.46/good/DHCPv6PrefixDelegation_Assign_1.target
[Unit]
Description=Foo
[DHCPv6PrefixDelegation]
Assign=1

