# tests/auto/InvalidValue/2.46/good/EnhancedTransmissionSelection_Handle_0xffff.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
Handle=0xffff

