# tests/auto/InvalidValue/2.46/good/EnhancedTransmissionSelection_Parent_root.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
Parent=root

