# tests/auto/InvalidValue/2.46/good/EnhancedTransmissionSelection_StrictBands_1.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
StrictBands=1

