# tests/auto/InvalidValue/2.46/good/EnhancedTransmissionSelection_QuantumBytes_0K.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
QuantumBytes=0K

