# tests/auto/InvalidValue/2.46/good/EnhancedTransmissionSelection_QuantumBytes_0M.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
QuantumBytes=0M

