# tests/auto/InvalidValue/2.46/good/EnhancedTransmissionSelection_Parent_ingress.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
Parent=ingress

