# tests/auto/InvalidValue/2.46/good/EnhancedTransmissionSelection_StrictBands_16.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
StrictBands=16

