# tests/auto/InvalidValue/2.46/good/EnhancedTransmissionSelection_Bands_16.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
Bands=16

