# tests/auto/InvalidValue/2.46/good/GenericRandomEarlyDetection_Parent_root.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
Parent=root

