# tests/auto/InvalidValue/2.46/good/HeavyHitterFilter_Parent_abc.target
[Unit]
Description=Foo
[HeavyHitterFilter]
Parent=abc

