# tests/auto/InvalidValue/2.46/good/HeavyHitterFilter_Handle_0xffff.target
[Unit]
Description=Foo
[HeavyHitterFilter]
Handle=0xffff

