# tests/auto/InvalidValue/2.46/good/HeavyHitterFilter_Handle_0x1.target
[Unit]
Description=Foo
[HeavyHitterFilter]
Handle=0x1

