# tests/auto/InvalidValue/2.46/good/HierarchyTokenBucketClass_BufferBytes_9999999999999K.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
BufferBytes=9999999999999K

