# tests/auto/InvalidValue/2.46/good/HierarchyTokenBucketClass_ClassId_abc.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
ClassId=abc

