# tests/auto/InvalidValue/2.46/good/HierarchyTokenBucketClass_CeilRate_9999999999999K.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
CeilRate=9999999999999K

