# tests/auto/InvalidValue/2.46/good/HierarchyTokenBucketClass_OverheadBytes_0G.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
OverheadBytes=0G

