# tests/auto/InvalidValue/2.46/good/HierarchyTokenBucket_Handle_0x1.target
[Unit]
Description=Foo
[HierarchyTokenBucket]
Handle=0x1

