# tests/auto/InvalidValue/2.46/good/HierarchyTokenBucketClass_Parent_root.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
Parent=root

