# tests/auto/InvalidValue/2.46/good/HierarchyTokenBucket_Parent_abc.target
[Unit]
Description=Foo
[HierarchyTokenBucket]
Parent=abc

