# tests/auto/InvalidValue/2.46/good/HierarchyTokenBucket_Parent_root.target
[Unit]
Description=Foo
[HierarchyTokenBucket]
Parent=root

