# tests/auto/InvalidValue/2.46/good/HierarchyTokenBucketClass_CeilRate_0M.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
CeilRate=0M

