# tests/auto/InvalidValue/2.46/good/HierarchyTokenBucketClass_CeilBufferBytes_9999999999999K.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
CeilBufferBytes=9999999999999K

