# tests/auto/InvalidValue/2.46/good/HierarchyTokenBucketClass_Priority_9999999999999.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
Priority=9999999999999

