# tests/auto/InvalidValue/2.46/good/HierarchyTokenBucket_Parent_clsact.target
[Unit]
Description=Foo
[HierarchyTokenBucket]
Parent=clsact

