# tests/auto/InvalidValue/2.46/good/HierarchyTokenBucketClass_Rate_9999999999999K.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
Rate=9999999999999K

