# tests/auto/InvalidValue/2.46/good/HierarchyTokenBucketClass_MTUBytes_9999999999999M.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
MTUBytes=9999999999999M

