# tests/auto/InvalidValue/2.46/good/HierarchyTokenBucketClass_QuantumBytes_9999999999999M.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
QuantumBytes=9999999999999M

