# tests/auto/InvalidValue/2.46/good/IPv6AcceptRA_DHCPv6Client_on.target
[Unit]
Description=Foo
[IPv6AcceptRA]
DHCPv6Client=on

