# tests/auto/InvalidValue/2.46/good/Link_AutoNegotiationFlowControl_0.target
[Unit]
Description=Foo
[Link]
AutoNegotiationFlowControl=0

