# tests/auto/InvalidValue/2.46/good/Link_AutoNegotiationFlowControl_on.target
[Unit]
Description=Foo
[Link]
AutoNegotiationFlowControl=on

