# tests/auto/InvalidValue/2.46/good/Link_TxFlowControl_on.target
[Unit]
Description=Foo
[Link]
TxFlowControl=on

