# tests/auto/InvalidValue/2.46/good/Link_RxFlowControl_no.target
[Unit]
Description=Foo
[Link]
RxFlowControl=no

