# tests/auto/InvalidValue/2.46/good/Link_RxFlowControl_false.target
[Unit]
Description=Foo
[Link]
RxFlowControl=false

