# tests/auto/InvalidValue/2.46/good/Link_TxFlowControl_off.target
[Unit]
Description=Foo
[Link]
TxFlowControl=off

