# tests/auto/InvalidValue/2.46/good/Link_AutoNegotiationFlowControl_false.target
[Unit]
Description=Foo
[Link]
AutoNegotiationFlowControl=false

