# tests/auto/InvalidValue/2.46/good/Link_RxFlowControl_0.target
[Unit]
Description=Foo
[Link]
RxFlowControl=0

