# tests/auto/InvalidValue/2.46/good/Link_AutoNegotiationFlowControl_true.target
[Unit]
Description=Foo
[Link]
AutoNegotiationFlowControl=true

