# tests/auto/InvalidValue/2.46/good/Link_AutoNegotiationFlowControl_off.target
[Unit]
Description=Foo
[Link]
AutoNegotiationFlowControl=off

