# tests/auto/InvalidValue/2.46/good/Link_RxFlowControl_off.target
[Unit]
Description=Foo
[Link]
RxFlowControl=off

