# tests/auto/InvalidValue/2.46/good/LLDP_MUDURL_abc.target
[Unit]
Description=Foo
[LLDP]
MUDURL=abc

