# tests/auto/InvalidValue/2.46/good/Network_IPv6LinkLocalAddressGenerationMode_none.target
[Unit]
Description=Foo
[Network]
IPv6LinkLocalAddressGenerationMode=none

