# tests/auto/InvalidValue/2.46/good/PFIFO_Parent_ingress.target
[Unit]
Description=Foo
[PFIFO]
Parent=ingress

