# tests/auto/InvalidValue/2.46/good/PFIFOHeadDrop_PacketLimit_0.target
[Unit]
Description=Foo
[PFIFOHeadDrop]
PacketLimit=0

