# tests/auto/InvalidValue/2.46/good/PFIFO_PacketLimit_4294967294.target
[Unit]
Description=Foo
[PFIFO]
PacketLimit=4294967294

