# tests/auto/InvalidValue/2.46/good/PFIFO_PacketLimit_0.target
[Unit]
Description=Foo
[PFIFO]
PacketLimit=0

