# tests/auto/InvalidValue/2.46/good/PFIFOFast_Parent_ingress.target
[Unit]
Description=Foo
[PFIFOFast]
Parent=ingress

