# tests/auto/InvalidValue/2.46/good/PIE_Parent_abc.target
[Unit]
Description=Foo
[PIE]
Parent=abc

