# tests/auto/InvalidValue/2.46/good/PIE_Parent_ingress.target
[Unit]
Description=Foo
[PIE]
Parent=ingress

