# tests/auto/InvalidValue/2.46/good/PIE_Parent_clsact.target
[Unit]
Description=Foo
[PIE]
Parent=clsact

