# tests/auto/InvalidValue/2.46/good/PIE_Handle_0xffff.target
[Unit]
Description=Foo
[PIE]
Handle=0xffff

