# tests/auto/InvalidValue/2.46/good/QuickFairQueueingClass_Weight_1.target
[Unit]
Description=Foo
[QuickFairQueueingClass]
Weight=1

