# tests/auto/InvalidValue/2.46/good/QuickFairQueueingClass_Parent_root.target
[Unit]
Description=Foo
[QuickFairQueueingClass]
Parent=root

