# tests/auto/InvalidValue/2.46/good/QuickFairQueueing_Handle_0x1.target
[Unit]
Description=Foo
[QuickFairQueueing]
Handle=0x1

