# tests/auto/InvalidValue/2.46/good/Service_TimeoutStopFailureMode_kill.target
[Unit]
Description=Foo
[Service]
TimeoutStopFailureMode=kill

