# tests/auto/InvalidValue/2.46/good/Service_TimeoutStartFailureMode_kill.target
[Unit]
Description=Foo
[Service]
TimeoutStartFailureMode=kill

