# tests/auto/InvalidValue/2.46/good/Service_RootHashSignature_abc.target
[Unit]
Description=Foo
[Service]
RootHashSignature=abc

