# tests/auto/InvalidValue/2.46/good/Service_CoredumpFilter_shared-file-backed.target
[Unit]
Description=Foo
[Service]
CoredumpFilter=shared-file-backed

