# tests/auto/InvalidValue/2.46/good/Service_CoredumpFilter_elf-headers.target
[Unit]
Description=Foo
[Service]
CoredumpFilter=elf-headers

