# tests/auto/InvalidValue/2.46/good/Socket_CoredumpFilter_default.target
[Unit]
Description=Foo
[Socket]
CoredumpFilter=default

