# tests/auto/InvalidValue/2.46/good/Socket_CoredumpFilter_shared-dax.target
[Unit]
Description=Foo
[Socket]
CoredumpFilter=shared-dax

