# tests/auto/InvalidValue/2.46/good/Socket_PassPacketInfo_false.target
[Unit]
Description=Foo
[Socket]
PassPacketInfo=false

