# tests/auto/InvalidValue/2.46/good/Socket_PassPacketInfo_0.target
[Unit]
Description=Foo
[Socket]
PassPacketInfo=0

