# tests/auto/InvalidValue/2.46/good/Socket_PassPacketInfo_on.target
[Unit]
Description=Foo
[Socket]
PassPacketInfo=on

