# tests/auto/InvalidValue/2.46/good/Socket_CoredumpFilter_elf-headers.target
[Unit]
Description=Foo
[Socket]
CoredumpFilter=elf-headers

