# tests/auto/InvalidValue/2.46/good/Socket_CoredumpFilter_all.target
[Unit]
Description=Foo
[Socket]
CoredumpFilter=all

