# tests/auto/InvalidValue/2.46/good/Socket_CoredumpFilter_shared-huge.target
[Unit]
Description=Foo
[Socket]
CoredumpFilter=shared-huge

