# tests/auto/InvalidValue/2.46/good/Socket_RootHashSignature_abc.target
[Unit]
Description=Foo
[Socket]
RootHashSignature=abc

