# tests/auto/InvalidValue/2.46/good/Socket_PassPacketInfo_off.target
[Unit]
Description=Foo
[Socket]
PassPacketInfo=off

