# tests/auto/InvalidValue/2.46/good/SR-IOV_VLANId_4096.target
[Unit]
Description=Foo
[SR-IOV]
VLANId=4096

