# tests/auto/InvalidValue/2.46/good/SR-IOV_VirtualFunction_2147483646.target
[Unit]
Description=Foo
[SR-IOV]
VirtualFunction=2147483646

