# tests/auto/InvalidValue/2.46/good/SR-IOV_Trust_false.target
[Unit]
Description=Foo
[SR-IOV]
Trust=false

