# tests/auto/InvalidValue/2.46/good/SR-IOV_LinkState_auto.target
[Unit]
Description=Foo
[SR-IOV]
LinkState=auto

