# tests/auto/InvalidValue/2.46/good/SR-IOV_QueryReceiveSideScaling_1.target
[Unit]
Description=Foo
[SR-IOV]
QueryReceiveSideScaling=1

