# tests/auto/InvalidValue/2.46/good/SR-IOV_QueryReceiveSideScaling_yes.target
[Unit]
Description=Foo
[SR-IOV]
QueryReceiveSideScaling=yes

