# tests/auto/InvalidValue/2.46/good/SR-IOV_VLANProtocol_802.1ad.target
[Unit]
Description=Foo
[SR-IOV]
VLANProtocol=802.1ad

