# tests/auto/InvalidValue/2.46/good/SR-IOV_MACSpoofCheck_on.target
[Unit]
Description=Foo
[SR-IOV]
MACSpoofCheck=on

