# tests/auto/InvalidValue/2.46/good/SR-IOV_Trust_0.target
[Unit]
Description=Foo
[SR-IOV]
Trust=0

