# tests/auto/InvalidValue/2.46/good/StochasticFairBlue_Handle_0x1.target
[Unit]
Description=Foo
[StochasticFairBlue]
Handle=0x1

