# tests/auto/InvalidValue/2.46/good/StochasticFairBlue_Parent_root.target
[Unit]
Description=Foo
[StochasticFairBlue]
Parent=root

