# tests/auto/InvalidValue/2.46/good/StochasticFairBlue_Parent_clsact.target
[Unit]
Description=Foo
[StochasticFairBlue]
Parent=clsact

