# tests/auto/InvalidValue/2.46/good/StochasticFairBlue_Parent_ingress.target
[Unit]
Description=Foo
[StochasticFairBlue]
Parent=ingress

