# tests/auto/InvalidValue/2.46/good/StochasticFairBlue_PacketLimit_0.target
[Unit]
Description=Foo
[StochasticFairBlue]
PacketLimit=0

