# tests/auto/InvalidValue/2.46/good/Swap_CoredumpFilter_16.target
[Unit]
Description=Foo
[Swap]
CoredumpFilter=16

