# tests/auto/InvalidValue/2.46/good/Swap_CoredumpFilter_shared-dax.target
[Unit]
Description=Foo
[Swap]
CoredumpFilter=shared-dax

