# tests/auto/InvalidValue/2.46/good/Swap_CoredumpFilter_all.target
[Unit]
Description=Foo
[Swap]
CoredumpFilter=all

