# tests/auto/InvalidValue/2.46/good/Swap_RootHashSignature_abc.target
[Unit]
Description=Foo
[Swap]
RootHashSignature=abc

