# tests/auto/InvalidValue/2.46/good/Swap_CoredumpFilter_elf-headers.target
[Unit]
Description=Foo
[Swap]
CoredumpFilter=elf-headers

