# tests/auto/InvalidValue/2.46/good/Swap_CoredumpFilter_private-anonymous.target
[Unit]
Description=Foo
[Swap]
CoredumpFilter=private-anonymous

