# tests/auto/InvalidValue/2.46/good/Swap_CoredumpFilter_shared-file-backed.target
[Unit]
Description=Foo
[Swap]
CoredumpFilter=shared-file-backed

