# tests/auto/InvalidValue/2.46/good/Unit_AssertPathIsEncrypted_--dev--.target
[Unit]
Description=Foo
Description=abc
AssertPathIsEncrypted=/dev/

