# tests/auto/InvalidValue/2.46/good/Unit_ConditionPathIsEncrypted_--dev--.target
[Unit]
Description=Foo
Description=abc
ConditionPathIsEncrypted=/dev/

