# systemdlint/tests/auto/InvalidValue/2.42/good/Mount_CPUQuotaPeriodSec_21us.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
CPUQuotaPeriodSec=21us

