# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_1day.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=1day

