# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_28h.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=28h

