# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_19d.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=19d

