# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_10day.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=10day

