# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_0sec.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=0sec

