# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_24years.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=24years

