# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_20usec.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=20usec

