# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_23week.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=23week

