# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_13minute.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=13minute

