# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_9month.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=9month

