# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_20w.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=20w

