# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_8minute.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=8minute

