# systemdlint/tests/auto/InvalidValue/2.42/good/Service_ProtectHostname_on.target
[Unit]
Description=Foo
[Service]
ProtectHostname=on

