# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_17minutes.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=17minutes

