# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_7seconds.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=7seconds

