# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_26months.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=26months

