# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_27msec.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=27msec

