# systemdlint/tests/auto/InvalidValue/2.42/good/Service_ProtectHostname_off.target
[Unit]
Description=Foo
[Service]
ProtectHostname=off

