# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_22weeks.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=22weeks

