# systemdlint/tests/auto/InvalidValue/2.42/good/Service_ProtectHostname_1.target
[Unit]
Description=Foo
[Service]
ProtectHostname=1

