# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_4ms.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=4ms

