# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_15ms.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=15ms

