# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_12second.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=12second

