# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_25s.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=25s

