# systemdlint/tests/auto/InvalidValue/2.42/good/Service_ProtectHostname_0.target
[Unit]
Description=Foo
[Service]
ProtectHostname=0

