# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_16seconds.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=16seconds

