# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_3hours.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=3hours

