# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_5hour.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=5hour

