# systemdlint/tests/auto/InvalidValue/2.42/good/Service_ProtectHostname_true.target
[Unit]
Description=Foo
[Service]
ProtectHostname=true

