# systemdlint/tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_21usec.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=21usec

