# systemdlint/tests/auto/InvalidValue/2.42/good/Socket_CPUQuotaPeriodSec_6days.target
[Unit]
Description=Foo
[Socket]
CPUQuotaPeriodSec=6days

