# systemdlint/tests/auto/InvalidValue/2.42/good/Socket_CPUQuotaPeriodSec_22weeks.target
[Unit]
Description=Foo
[Socket]
CPUQuotaPeriodSec=22weeks

