# systemdlint/tests/auto/InvalidValue/2.42/good/Socket_CPUQuotaPeriodSec_26months.target
[Unit]
Description=Foo
[Socket]
CPUQuotaPeriodSec=26months

