# systemdlint/tests/auto/InvalidValue/2.42/good/Socket_CPUQuotaPeriodSec_0sec.target
[Unit]
Description=Foo
[Socket]
CPUQuotaPeriodSec=0sec

