# systemdlint/tests/auto/InvalidValue/2.42/good/Socket_CPUQuotaPeriodSec_13minute.target
[Unit]
Description=Foo
[Socket]
CPUQuotaPeriodSec=13minute

