# systemdlint/tests/auto/InvalidValue/2.42/good/Socket_CPUQuotaPeriodSec_14M.target
[Unit]
Description=Foo
[Socket]
CPUQuotaPeriodSec=14M

