# systemdlint/tests/auto/InvalidValue/2.42/good/Socket_CPUQuotaPeriodSec_25s.target
[Unit]
Description=Foo
[Socket]
CPUQuotaPeriodSec=25s

