# systemdlint/tests/auto/InvalidValue/2.42/good/Socket_ProtectHostname_true.target
[Unit]
Description=Foo
[Socket]
ProtectHostname=true

