# systemdlint/tests/auto/InvalidValue/2.42/good/Socket_ProtectHostname_no.target
[Unit]
Description=Foo
[Socket]
ProtectHostname=no

