# systemdlint/tests/auto/InvalidValue/2.42/good/Socket_CPUQuotaPeriodSec_17minutes.target
[Unit]
Description=Foo
[Socket]
CPUQuotaPeriodSec=17minutes

