# systemdlint/tests/auto/InvalidValue/2.42/good/Socket_CPUQuotaPeriodSec_11y.target
[Unit]
Description=Foo
[Socket]
CPUQuotaPeriodSec=11y

