# systemdlint/tests/auto/InvalidValue/2.42/good/Socket_CPUQuotaPeriodSec_10day.target
[Unit]
Description=Foo
[Socket]
CPUQuotaPeriodSec=10day

