# systemdlint/tests/auto/InvalidValue/2.42/good/Socket_CPUQuotaPeriodSec_5hour.target
[Unit]
Description=Foo
[Socket]
CPUQuotaPeriodSec=5hour

