# systemdlint/tests/auto/InvalidValue/2.42/good/Socket_ProtectHostname_on.target
[Unit]
Description=Foo
[Socket]
ProtectHostname=on

