# systemdlint/tests/auto/InvalidValue/2.42/good/Swap_CPUQuotaPeriodSec_0sec.target
[Unit]
Description=Foo
[Swap]
CPUQuotaPeriodSec=0sec

