# systemdlint/tests/auto/InvalidValue/2.42/good/Swap_ProtectHostname_false.target
[Unit]
Description=Foo
[Swap]
ProtectHostname=false

