# systemdlint/tests/auto/InvalidValue/2.42/good/Swap_CPUQuotaPeriodSec_17second.target
[Unit]
Description=Foo
[Swap]
CPUQuotaPeriodSec=17second

