# systemdlint/tests/auto/InvalidValue/2.42/good/Swap_CPUQuotaPeriodSec_7seconds.target
[Unit]
Description=Foo
[Swap]
CPUQuotaPeriodSec=7seconds

