# systemdlint/tests/auto/InvalidValue/2.42/good/Swap_CPUQuotaPeriodSec_2d.target
[Unit]
Description=Foo
[Swap]
CPUQuotaPeriodSec=2d

