# systemdlint/tests/auto/InvalidValue/2.42/good/Swap_CPUQuotaPeriodSec_18year.target
[Unit]
Description=Foo
[Swap]
CPUQuotaPeriodSec=18year

