# systemdlint/tests/auto/InvalidValue/2.42/good/Swap_ProtectHostname_on.target
[Unit]
Description=Foo
[Swap]
ProtectHostname=on

