# systemdlint/tests/auto/InvalidValue/2.42/good/Swap_CPUQuotaPeriodSec_15min.target
[Unit]
Description=Foo
[Swap]
CPUQuotaPeriodSec=15min

