# systemdlint/tests/auto/InvalidValue/2.42/good/Swap_CPUQuotaPeriodSec_14msec.target
[Unit]
Description=Foo
[Swap]
CPUQuotaPeriodSec=14msec

