# systemdlint/tests/auto/InvalidValue/2.42/good/Swap_CPUQuotaPeriodSec_6days.target
[Unit]
Description=Foo
[Swap]
CPUQuotaPeriodSec=6days

