# systemdlint/tests/auto/InvalidValue/2.42/good/Swap_CPUQuotaPeriodSec_11y.target
[Unit]
Description=Foo
[Swap]
CPUQuotaPeriodSec=11y

