# systemdlint/tests/auto/InvalidValue/2.42/good/Swap_CPUQuotaPeriodSec_19s.target
[Unit]
Description=Foo
[Swap]
CPUQuotaPeriodSec=19s

