# systemdlint/tests/auto/InvalidValue/2.42/good/Swap_CPUQuotaPeriodSec_infinity.target
[Unit]
Description=Foo
[Swap]
CPUQuotaPeriodSec=infinity

