# systemdlint/tests/auto/InvalidValue/2.42/good/Swap_ProtectHostname_yes.target
[Unit]
Description=Foo
[Swap]
ProtectHostname=yes

