# systemdlint/tests/auto/InvalidValue/2.42/good/Unit_ConditionCPUs_0.target
[Unit]
Description=Foo
Description=abc
ConditionCPUs=0

