# systemdlint/tests/auto/InvalidNumericBase/2.31/good/CAN_Bitrate_3000.target
[Unit]
Description=Foo
[CAN]
Bitrate=3000

