# systemdlint/tests/auto/InvalidNumericBase/2.31/good/Link_MTUBytes_0K.target
[Unit]
Description=Foo
[Link]
MACAddress=abc
MTUBytes=0K

