# systemdlint/tests/auto/InvalidNumericBase/2.31/good/Login_RuntimeDirectorySize_6144.target
[Unit]
Description=Foo
[Login]
RuntimeDirectorySize=6144

