# systemdlint/tests/auto/InvalidNumericBase/2.31/good/Login_RuntimeDirectorySize_3072.target
[Unit]
Description=Foo
[Login]
RuntimeDirectorySize=3072

