# systemdlint/tests/auto/InvalidNumericBase/2.31/good/Login_RuntimeDirectorySize_0K.target
[Unit]
Description=Foo
[Login]
RuntimeDirectorySize=0K

