# systemdlint/tests/auto/InvalidNumericBase/2.31/good/Login_RuntimeDirectorySize_5120.target
[Unit]
Description=Foo
[Login]
RuntimeDirectorySize=5120

