# systemdlint/tests/auto/InvalidNumericBase/2.31/good/Login_RuntimeDirectorySize_8192.target
[Unit]
Description=Foo
[Login]
RuntimeDirectorySize=8192

