# systemdlint/tests/auto/InvalidNumericBase/2.31/good/Route_MTUBytes_3072.target
[Unit]
Description=Foo
[Route]
MTUBytes=3072

