# systemdlint/tests/auto/InvalidNumericBase/2.31/good/Route_MTUBytes_8192.target
[Unit]
Description=Foo
[Route]
MTUBytes=8192

