# tests/auto/OptionTooNew/2.45/good/FairQueueingControlledDelay_IntervalSec_0days.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=123

