# tests/auto/OptionTooNew/2.45/good/FairQueueing_CEThresholdSec_0days.target
[Unit]
Description=Foo
[FairQueueing]
CEThresholdSec=123

