# tests/auto/OptionTooNew/2.45/good/Link_AlternativeNamesPolicy_database.target
[Unit]
Description=Foo
[Link]
AlternativeNamesPolicy=123

