# tests/auto/OptionTooNew/2.45/good/NetworkEmulator_Parent_root.target
[Unit]
Description=Foo
[NetworkEmulator]
Parent=123

