# tests/auto/OptionTooNew/2.45/good/NetworkEmulator_PacketLimit_0.target
[Unit]
Description=Foo
[NetworkEmulator]
PacketLimit=123

