# tests/auto/OptionTooNew/2.45/good/Swap_ProtectClock_yes.target
[Unit]
Description=Foo
[Swap]
ProtectClock=123

