# tests/auto/OptionTooNew/2.45/good/TokenBucketFilter_LimitSize_0K.target
[Unit]
Description=Foo
[TokenBucketFilter]
LimitSize=123

