# systemdlint/tests/auto/OptionTooNew/2.28/bad/Service_RuntimeMaxSec_0sec.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=123

