# tests/auto/OptionTooNew/2.37/bad/Mount_TemporaryFileSystem_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
TemporaryFileSystem=123

