# tests/auto/OptionTooNew/2.37/bad/Service_TemporaryFileSystem_abc.target
[Unit]
Description=Foo
[Service]
TemporaryFileSystem=123

