# tests/auto/OptionTooNew/2.37/bad/Unit_ConditionKernelVersion_abc.target
[Unit]
Description=Foo
Description=abc
ConditionKernelVersion=123

