# systemdlint/tests/auto/OptionTooNew/2.28/good/Service_PassEnvironment_abc.target
[Unit]
Description=Foo
[Service]
PassEnvironment=123

