# tests/auto/OptionTooNew/2.48/bad/DHCPServerStaticLease_Address_1.1.1.1--0.target
[Unit]
Description=Foo
[DHCPServerStaticLease]
Address=123

