# tests/auto/OptionTooNew/2.48/bad/DHCPServer_ServerAddress_1.1.1.1--0.target
[Unit]
Description=Foo
[DHCPServer]
ServerAddress=123

