# tests/auto/OptionTooNew/2.48/bad/Link_RequiredFamilyForOnline_ipv4.target
[Unit]
Description=Foo
[Link]
RequiredFamilyForOnline=123

