# tests/auto/OptionTooNew/2.48/bad/Scope_SocketBindAllow_abc.target
[Unit]
Description=Foo
[Scope]
SocketBindAllow=123

