# tests/auto/OptionTooNew/2.48/bad/Service_SocketBindDeny_abc.target
[Unit]
Description=Foo
[Service]
SocketBindDeny=123

