# tests/auto/OptionTooNew/2.48/bad/Socket_SocketBindAllow_abc.target
[Unit]
Description=Foo
[Socket]
SocketBindAllow=123

