# tests/auto/OptionTooNew/2.48/bad/Swap_SocketBindDeny_abc.target
[Unit]
Description=Foo
[Swap]
SocketBindDeny=123

