# tests/auto/OptionTooNew/2.48/bad/Unit_StopPropagatedFrom_syslog.service.target
[Unit]
Description=Foo
Description=abc
StopPropagatedFrom=123

