# tests/auto/OptionTooNew/2.48/bad/Unit_ConditionFirmware_abc.target
[Unit]
Description=Foo
Description=abc
ConditionFirmware=123

