# tests/auto/OptionTooNew/2.48/bad/Unit_Upholds_syslog.service.target
[Unit]
Description=Foo
Description=abc
Upholds=123

