# tests/auto/OptionTooNew/2.48/bad/Upload_ServerKeyFile_abc.target
[Unit]
Description=Foo
[Upload]
ServerKeyFile=123

