# tests/auto/OptionDeprecated/2.50/good/DHCPv6PrefixDelegation_SubnetId_0x0.target
[Unit]
Description=Foo
[DHCPv6PrefixDelegation]
SubnetId=123

