# tests/auto/InvalidValue/2.33/good/Socket_BindPaths_abc.target
[Unit]
Description=Foo
[Socket]
BindPaths=abc

