# tests/auto/InvalidValue/2.33/good/Swap_BindReadOnlyPaths_abc.target
[Unit]
Description=Foo
[Swap]
BindReadOnlyPaths=abc

