# tests/auto/OptionTooNew/2.15/good/Unit_ConditionNeedsUpdate_abc.target
[Unit]
Description=Foo
Description=abc
ConditionNeedsUpdate=123

