# systemdlint/tests/auto/SettingRestricted/2.31/bad/DHCP_ClientIdentifier_macAnonymize%3D1.target
[Unit]
Description=Foo
[DHCP]
ClientIdentifier=mac
Anonymize=1
