# systemdlint/tests/auto/SettingRestricted/2.31/bad/DHCP_UseRoutes_yesAnonymize%3D1.target
[Unit]
Description=Foo
[DHCP]
UseRoutes=yes
Anonymize=1
