# systemdlint/tests/auto/SettingRestricted/2.31/bad/DHCP_VendorClassIdentifier_abcAnonymize%3Dtrue.target
[Unit]
Description=Foo
[DHCP]
VendorClassIdentifier=abc
Anonymize=true
