# systemdlint/tests/auto/SettingRestricted/2.31/bad/DHCP_UseTimezone_yesAnonymize%3D1.target
[Unit]
Description=Foo
[DHCP]
UseTimezone=yes
Anonymize=1
