# systemdlint/tests/auto/SettingRestricted/2.31/bad/DHCP_UseMTU_yesAnonymize%3Dyes-MTUBytes%3DNone.target
[Unit]
Description=Foo
[DHCP]
UseMTU=yes
Anonymize=yes
MTUBytes=None
