# systemdlint/tests/auto/SettingRestricted/2.31/bad/DHCP_ClientIdentifier_macAnonymize%3Dyes.target
[Unit]
Description=Foo
[DHCP]
ClientIdentifier=mac
Anonymize=yes
