# systemdlint/tests/auto/SettingRestricted/2.31/bad/DHCP_UseMTU_yesAnonymize%3D1-MTUBytes%3DNone.target
[Unit]
Description=Foo
[DHCP]
UseMTU=yes
Anonymize=1
MTUBytes=None
