# systemdlint/tests/auto/SettingRestricted/2.31/bad/Network_PrimarySlave_yesMode%3Dactive-backup.target
[Unit]
Description=Foo
[Network]
PrimarySlave=yes
Mode=active-backup
