# systemdlint/tests/auto/SettingRestricted/2.31/bad/Network_PrimarySlave_yesMode%3Dbalance-tlb.target
[Unit]
Description=Foo
[Network]
PrimarySlave=yes
Mode=balance-tlb
