# systemdlint/tests/auto/SettingRestricted/2.31/bad/Service_Restart_noType%3Doneshot.target
[Unit]
Description=Foo
[Service]
Restart=no
Type=oneshot
