# systemdlint/tests/auto/OptionTooNew/2.31/good/Mount_InaccessiblePaths_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
InaccessiblePaths=123

