# systemdlint/tests/auto/OptionTooNew/2.31/good/Service_ReadWriteDirectories_abc.target
[Unit]
Description=Foo
[Service]
ReadWriteDirectories=123

