# systemdlint/tests/auto/OptionTooNew/2.31/good/Socket_RestrictRealtime_yes.target
[Unit]
Description=Foo
[Socket]
RestrictRealtime=123

