# systemdlint/tests/auto/OptionTooNew/2.31/good/Socket_InaccessiblePaths_abc.target
[Unit]
Description=Foo
[Socket]
InaccessiblePaths=123

