# systemdlint/tests/auto/OptionTooNew/2.31/good/Socket_ReadOnlyPaths_abc.target
[Unit]
Description=Foo
[Socket]
ReadOnlyPaths=123

