# systemdlint/tests/auto/OptionTooNew/2.31/good/Swap_ReadWritePaths_abc.target
[Unit]
Description=Foo
[Swap]
ReadWritePaths=123

