# systemdlint/tests/auto/OptionTooNew/2.31/good/Swap_ReadOnlyDirectories_abc.target
[Unit]
Description=Foo
[Swap]
ReadOnlyDirectories=123

