# systemdlint/tests/auto/OptionTooNew/2.31/good/Swap_InaccessableDirectories_abc.target
[Unit]
Description=Foo
[Swap]
InaccessableDirectories=123

