# systemdlint/tests/auto/OptionTooNew/2.31/good/Swap_ReadWriteDirectories_abc.target
[Unit]
Description=Foo
[Swap]
ReadWriteDirectories=123

