# tests/auto/OptionTooNew/2.38/good/Unit_ConditionControlGroupController_abc.target
[Unit]
Description=Foo
Description=abc
ConditionControlGroupController=123

