# tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_24w.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=24w

