# tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_14msec.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=14msec

