# tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_2d.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=2d

