# tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_4hour.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=4hour

