# tests/auto/InvalidValue/2.42/good/Service_ProtectHostname_false.target
[Unit]
Description=Foo
[Service]
ProtectHostname=false

