# tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_17second.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=17second

