# tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_infinity.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=infinity

