# tests/auto/InvalidValue/2.42/good/Service_NetworkNamespacePath_--dev--.target
[Unit]
Description=Foo
[Service]
NetworkNamespacePath=/dev/

