# tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_6h.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=6h

