# tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_11m.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=11m

