# tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_18sec.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=18sec

