# tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_12month.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=12month

