# tests/auto/InvalidValue/2.42/good/Service_ProtectHostname_yes.target
[Unit]
Description=Foo
[Service]
ProtectHostname=yes

