# tests/auto/InvalidValue/2.42/good/Service_CPUQuotaPeriodSec_7minutes.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=7minutes

