# tests/auto/InvalidValue/2.42/good/Socket_ProtectHostname_1.target
[Unit]
Description=Foo
[Socket]
ProtectHostname=1

