# tests/auto/InvalidValue/2.42/good/Socket_ProtectHostname_off.target
[Unit]
Description=Foo
[Socket]
ProtectHostname=off

