# tests/auto/InvalidValue/2.42/good/Socket_ProtectHostname_yes.target
[Unit]
Description=Foo
[Socket]
ProtectHostname=yes

