# tests/auto/InvalidValue/2.42/good/Socket_NetworkNamespacePath_--dev--.target
[Unit]
Description=Foo
[Socket]
NetworkNamespacePath=/dev/

