# tests/auto/InvalidValue/2.42/good/Socket_CPUQuotaPeriodSec_3hours.target
[Unit]
Description=Foo
[Socket]
CPUQuotaPeriodSec=3hours

