# tests/auto/InvalidValue/2.42/good/Socket_CPUQuotaPeriodSec_27y.target
[Unit]
Description=Foo
[Socket]
CPUQuotaPeriodSec=27y

