# tests/auto/InvalidValue/2.42/good/Socket_ProtectHostname_false.target
[Unit]
Description=Foo
[Socket]
ProtectHostname=false

