# tests/auto/InvalidValue/2.42/good/Swap_ProtectHostname_true.target
[Unit]
Description=Foo
[Swap]
ProtectHostname=true

