# tests/auto/InvalidValue/2.42/good/Swap_NetworkNamespacePath_--dev--.target
[Unit]
Description=Foo
[Swap]
NetworkNamespacePath=/dev/

