# tests/auto/InvalidValue/2.42/good/Swap_CPUQuotaPeriodSec_24w.target
[Unit]
Description=Foo
[Swap]
CPUQuotaPeriodSec=24w

