# tests/auto/InvalidValue/2.42/good/Swap_ProtectHostname_no.target
[Unit]
Description=Foo
[Swap]
ProtectHostname=no

