# tests/auto/InvalidValue/2.42/good/Swap_CPUQuotaPeriodSec_27y.target
[Unit]
Description=Foo
[Swap]
CPUQuotaPeriodSec=27y

