# tests/auto/InvalidValue/2.42/good/Swap_CPUQuotaPeriodSec_5hr.target
[Unit]
Description=Foo
[Swap]
CPUQuotaPeriodSec=5hr

