# tests/auto/InvalidValue/2.42/good/Unit_ConditionCPUs_9999999999999.target
[Unit]
Description=Foo
Description=abc
ConditionCPUs=9999999999999

