# systemdlint/tests/auto/SettingRequires/2.31/bad/Link_AllMulticast_yesMulticast%3D0.target
[Unit]
Description=Foo
[Link]
AllMulticast=yes

