# systemdlint/tests/auto/SettingRequires/2.31/bad/Login_IdleActionSec_0daysIdleAction%3DNone.target
[Unit]
Description=Foo
[Login]
IdleActionSec=0days

