# systemdlint/tests/auto/SettingRequires/2.31/bad/Service_BusName_abcType%3Ddbus.target
[Unit]
Description=Foo
[Service]
BusName=abc

