# tests/auto/InvalidValue/2.50/bad/CAKE_MPUBytes_257.target
[Unit]
Description=Foo
[CAKE]
MPUBytes=257

