# tests/auto/InvalidValue/2.50/bad/CAKE_FirewallMark_4294967296.target
[Unit]
Description=Foo
[CAKE]
FirewallMark=4294967296

