# tests/auto/InvalidValue/2.50/bad/CAKE_AutoRateIngress_baz.target
[Unit]
Description=Foo
[CAKE]
AutoRateIngress=baz

