# tests/auto/InvalidValue/2.50/bad/CAKE_UseRawPacketSize_foo.target
[Unit]
Description=Foo
[CAKE]
UseRawPacketSize=foo

