# tests/auto/InvalidValue/2.50/bad/CAKE_SplitGSO_baz.target
[Unit]
Description=Foo
[CAKE]
SplitGSO=baz

