# tests/auto/InvalidValue/2.50/bad/CAKE_UseRawPacketSize_baz.target
[Unit]
Description=Foo
[CAKE]
UseRawPacketSize=baz

