# tests/auto/InvalidValue/2.50/bad/CAKE_FirewallMark_0.target
[Unit]
Description=Foo
[CAKE]
FirewallMark=0

