# tests/auto/InvalidValue/2.50/bad/CAN_Loopback_foo.target
[Unit]
Description=Foo
[CAN]
Loopback=foo

