# tests/auto/InvalidValue/2.50/bad/CAN_Loopback_bar.target
[Unit]
Description=Foo
[CAN]
Loopback=bar

