# tests/auto/InvalidValue/2.50/bad/CAN_TimeQuantaNSec_-1.target
[Unit]
Description=Foo
[CAN]
TimeQuantaNSec=-1

