# tests/auto/InvalidValue/2.50/bad/CAN_DataPhaseBufferSegment1_-1.target
[Unit]
Description=Foo
[CAN]
DataPhaseBufferSegment1=-1

