# tests/auto/InvalidValue/2.50/bad/CAN_OneShot_bar.target
[Unit]
Description=Foo
[CAN]
OneShot=bar

