# tests/auto/InvalidValue/2.50/bad/CAN_DataPhaseBufferSegment2_-1.target
[Unit]
Description=Foo
[CAN]
DataPhaseBufferSegment2=-1

