# tests/auto/InvalidValue/2.50/bad/DHCPPrefixDelegation_RouteMetric_-1.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
RouteMetric=-1

