# tests/auto/InvalidValue/2.50/bad/DHCPv4_Use6RD_baz.target
[Unit]
Description=Foo
[DHCPv4]
Use6RD=baz

