# tests/auto/InvalidValue/2.50/bad/IPoIB_PartitionKey_65536.target
[Unit]
Description=Foo
[IPoIB]
PartitionKey=65536

