# tests/auto/InvalidValue/2.50/bad/IPoIB_IgnoreUserspaceMulticastGroup_foo.target
[Unit]
Description=Foo
[IPoIB]
IgnoreUserspaceMulticastGroup=foo

