# tests/auto/InvalidValue/2.50/bad/IPoIB_IgnoreUserspaceMulticastGroup_bar.target
[Unit]
Description=Foo
[IPoIB]
IgnoreUserspaceMulticastGroup=bar

