# tests/auto/InvalidValue/2.50/bad/Path_TriggerLimitBurst_10000000000000.target
[Unit]
Description=Foo
[Path]
TriggerLimitBurst=10000000000000

