# tests/auto/InvalidValue/2.50/bad/Path_TriggerLimitIntervalSec_0.5d.target
[Unit]
Description=Foo
[Path]
TriggerLimitIntervalSec=0.5d

