# tests/auto/InvalidValue/2.50/bad/RoutingPolicyRule_SuppressInterfaceGroup_-1.target
[Unit]
Description=Foo
[RoutingPolicyRule]
SuppressInterfaceGroup=-1

