# tests/auto/InvalidValue/2.50/bad/Scope_RuntimeRandomizedExtraSec_True.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=True

