# tests/auto/InvalidValue/2.50/bad/Scope_RuntimeRandomizedExtraSec_1p.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=1p

