# tests/auto/InvalidValue/2.50/bad/Service_TTYRows_10000000000000.target
[Unit]
Description=Foo
[Service]
TTYRows=10000000000000

