# tests/auto/InvalidValue/2.50/bad/Service_TTYRows_-1.target
[Unit]
Description=Foo
[Service]
TTYRows=-1

