# tests/auto/InvalidValue/2.50/bad/Service_TTYColumns_10000000000000.target
[Unit]
Description=Foo
[Service]
TTYColumns=10000000000000

