# tests/auto/InvalidValue/2.50/bad/Service_RuntimeRandomizedExtraSec_True.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=True

