# tests/auto/InvalidValue/2.50/bad/Socket_TTYColumns_10000000000000.target
[Unit]
Description=Foo
[Socket]
TTYColumns=10000000000000

