# tests/auto/InvalidValue/2.50/bad/Swap_TTYRows_-1.target
[Unit]
Description=Foo
[Swap]
TTYRows=-1

