# tests/auto/InvalidValue/2.50/bad/Swap_TTYRows_10000000000000.target
[Unit]
Description=Foo
[Swap]
TTYRows=10000000000000

