# tests/auto/InvalidValue/2.50/bad/Time_SaveIntervalSec_0.5d.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=0.5d

