# tests/auto/InvalidValue/2.50/bad/Time_SaveIntervalSec_3decades.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=3decades

