# tests/auto/SettingRequires/2.43/bad/DHCPv4_Anonymize_yesMACAddressPolicy%3D1.target
[Unit]
Description=Foo
[DHCPv4]
Anonymize=yes

