# systemdlint/tests/auto/InvalidValue/2.43/bad/BridgeFDB_Destination_0--abc.target
[Unit]
Description=Foo
[BridgeFDB]
MACAddress=01:23:45:67:89:ab
Destination=0/abc

