# systemdlint/tests/auto/InvalidValue/2.43/bad/Bridge_ProxyARPWiFi_foo.target
[Unit]
Description=Foo
[Bridge]
ProxyARPWiFi=foo

