# systemdlint/tests/auto/InvalidValue/2.43/bad/BridgeFDB_Destination_255.256.255.255.target
[Unit]
Description=Foo
[BridgeFDB]
MACAddress=abc
Destination=255.256.255.255

