# systemdlint/tests/auto/InvalidValue/2.43/bad/BridgeFDB_Destination_255.255.255.256.target
[Unit]
Description=Foo
[BridgeFDB]
MACAddress=abc
Destination=255.255.255.256

