# systemdlint/tests/auto/InvalidValue/2.43/bad/BridgeFDB_Destination_1.1.1.1---1.target
[Unit]
Description=Foo
[BridgeFDB]
MACAddress=abc
Destination=1.1.1.1/-1

