# systemdlint/tests/auto/InvalidValue/2.43/bad/BridgeFDB_Destination_abc.target
[Unit]
Description=Foo
[BridgeFDB]
MACAddress=abc
Destination=abc

