# systemdlint/tests/auto/InvalidValue/2.43/bad/BridgeFDB_AssociatedWith_1.target
[Unit]
Description=Foo
[BridgeFDB]
MACAddress=abc
AssociatedWith=1

