# systemdlint/tests/auto/InvalidValue/2.43/bad/Bridge_ProxyARP_baz.target
[Unit]
Description=Foo
[Bridge]
ProxyARP=baz

