# systemdlint/tests/auto/InvalidValue/2.43/bad/BridgeFDB_VNI_0.target
[Unit]
Description=Foo
[BridgeFDB]
MACAddress=abc
VNI=0

