# systemdlint/tests/auto/InvalidValue/2.43/bad/BridgeFDB_Destination_0.0.0..target
[Unit]
Description=Foo
[BridgeFDB]
MACAddress=abc
Destination=0.0.0.

