# systemdlint/tests/auto/InvalidValue/2.43/bad/DHCPv4_UseTimezone_foo.target
[Unit]
Description=Foo
[DHCPv4]
UseTimezone=foo

