# systemdlint/tests/auto/InvalidValue/2.43/bad/DHCPv6_UseRoutes_foo.target
[Unit]
Description=Foo
[DHCPv6]
UseRoutes=foo

