# systemdlint/tests/auto/InvalidValue/2.43/bad/DHCPv4_IAID_10000000000000.target
[Unit]
Description=Foo
[DHCPv4]
IAID=10000000000000

