# systemdlint/tests/auto/InvalidValue/2.43/bad/DHCPv4_UseHostname_foo.target
[Unit]
Description=Foo
[DHCPv4]
UseHostname=foo

