# systemdlint/tests/auto/InvalidValue/2.43/bad/DHCPv6_RouteMetric_10000000000000.target
[Unit]
Description=Foo
[DHCPv6]
RouteMetric=10000000000000

