# systemdlint/tests/auto/InvalidValue/2.43/bad/DHCPv4_RapidCommit_baz.target
[Unit]
Description=Foo
[DHCPv4]
RapidCommit=baz

