# systemdlint/tests/auto/InvalidValue/2.43/bad/DHCPv6_Anonymize_foo.target
[Unit]
Description=Foo
[DHCPv6]
Anonymize=foo

