# systemdlint/tests/auto/InvalidValue/2.43/bad/DHCPv4_RouteMetric_10000000000000.target
[Unit]
Description=Foo
[DHCPv4]
RouteMetric=10000000000000

