# systemdlint/tests/auto/InvalidValue/2.43/bad/DHCPv4_ForceDHCPv6PDOtherInformation_baz.target
[Unit]
Description=Foo
[DHCPv4]
ForceDHCPv6PDOtherInformation=baz

