# systemdlint/tests/auto/InvalidValue/2.43/bad/DHCPv6_UseDNS_foo.target
[Unit]
Description=Foo
[DHCPv6]
UseDNS=foo

