# systemdlint/tests/auto/InvalidValue/2.43/bad/DHCPv4_UseNTP_baz.target
[Unit]
Description=Foo
[DHCPv4]
UseNTP=baz

