# systemdlint/tests/auto/InvalidValue/2.43/bad/FooOverUDP_Peer_abc.target
[Unit]
Description=Foo
[FooOverUDP]
Peer=abc

