# systemdlint/tests/auto/InvalidValue/2.43/bad/FooOverUDP_Peer_0.0.0..target
[Unit]
Description=Foo
[FooOverUDP]
Peer=0.0.0.

