# systemdlint/tests/auto/InvalidValue/2.43/bad/FooOverUDP_Local_0--abc.target
[Unit]
Description=Foo
[FooOverUDP]
Local=0/abc

