# systemdlint/tests/auto/InvalidValue/2.43/bad/GENEVE_GenericProtocolExtension_foo.target
[Unit]
Description=Foo
[GENEVE]
GenericProtocolExtension=foo

