# systemdlint/tests/auto/InvalidValue/2.43/bad/GENEVE_IPDoNotFragment_bar.target
[Unit]
Description=Foo
[GENEVE]
IPDoNotFragment=bar

