# systemdlint/tests/auto/InvalidValue/2.43/bad/GENEVE_IPDoNotFragment_foo.target
[Unit]
Description=Foo
[GENEVE]
IPDoNotFragment=foo

