# systemdlint/tests/auto/InvalidValue/2.43/bad/Manager_NUMAPolicy_True.target
[Unit]
Description=Foo
[Manager]
NUMAPolicy=True

