# systemdlint/tests/auto/InvalidValue/2.43/bad/Manager_RebootWatchdogSec_True.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=True

