# systemdlint/tests/auto/InvalidValue/2.43/bad/Manager_RebootWatchdogSec_0.5d.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=0.5d

