# systemdlint/tests/auto/InvalidValue/2.43/bad/Manager_NUMAPolicy_1.target
[Unit]
Description=Foo
[Manager]
NUMAPolicy=1

