# systemdlint/tests/auto/InvalidValue/2.43/bad/Mount_DefaultMemoryMin_-1M.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
DefaultMemoryMin=-1M

