# systemdlint/tests/auto/InvalidValue/2.43/bad/Mount_NUMAPolicy_1.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
NUMAPolicy=1

