# systemdlint/tests/auto/InvalidValue/2.43/bad/Network_DefaultRouteOnDevice_baz.target
[Unit]
Description=Foo
[Network]
DefaultRouteOnDevice=baz

