# systemdlint/tests/auto/InvalidValue/2.43/bad/Route_FastOpenNoCookie_baz.target
[Unit]
Description=Foo
[Route]
FastOpenNoCookie=baz

