# systemdlint/tests/auto/InvalidValue/2.43/bad/Service_TimeoutAbortSec_True.target
[Unit]
Description=Foo
[Service]
TimeoutAbortSec=True

