# systemdlint/tests/auto/InvalidValue/2.43/bad/Service_NUMAPolicy_1.target
[Unit]
Description=Foo
[Service]
NUMAPolicy=1

