# systemdlint/tests/auto/InvalidValue/2.43/bad/Service_DefaultMemoryMin_-1K.target
[Unit]
Description=Foo
[Service]
DefaultMemoryMin=-1K

