# systemdlint/tests/auto/InvalidValue/2.43/bad/Service_NUMAPolicy_True.target
[Unit]
Description=Foo
[Service]
NUMAPolicy=True

