# systemdlint/tests/auto/InvalidValue/2.43/bad/Socket_NUMAPolicy_1.target
[Unit]
Description=Foo
[Socket]
NUMAPolicy=1

