# systemdlint/tests/auto/InvalidValue/2.43/bad/Socket_DefaultMemoryLow_-1T.target
[Unit]
Description=Foo
[Socket]
DefaultMemoryLow=-1T

