# systemdlint/tests/auto/InvalidValue/2.43/bad/Socket_DefaultMemoryLow_10000000000000B.target
[Unit]
Description=Foo
[Socket]
DefaultMemoryLow=10000000000000B

