# systemdlint/tests/auto/InvalidValue/2.43/bad/Socket_NUMAPolicy_True.target
[Unit]
Description=Foo
[Socket]
NUMAPolicy=True

