# systemdlint/tests/auto/InvalidValue/2.43/bad/Swap_NUMAPolicy_1.target
[Unit]
Description=Foo
[Swap]
NUMAPolicy=1

