# systemdlint/tests/auto/InvalidValue/2.43/bad/Swap_DefaultMemoryMin_-1T.target
[Unit]
Description=Foo
[Swap]
DefaultMemoryMin=-1T

