# systemdlint/tests/auto/InvalidValue/2.43/bad/Swap_NUMAPolicy_True.target
[Unit]
Description=Foo
[Swap]
NUMAPolicy=True

