# systemdlint/tests/auto/InvalidValue/2.43/bad/VXLAN_IPDoNotFragment_foo.target
[Unit]
Description=Foo
[VXLAN]
IPDoNotFragment=foo

