# systemdlint/tests/auto/InvalidValue/2.43/bad/VXLAN_IPDoNotFragment_baz.target
[Unit]
Description=Foo
[VXLAN]
IPDoNotFragment=baz

