# tests/auto/SettingRequires/2.31/good/DHCP_Anonymize_yesMACAddressPolicy%3Dyes.target
[Unit]
Description=Foo
[DHCP]
Anonymize=yes
MACAddressPolicy=yes
