# tests/auto/SettingRequires/2.31/good/Socket_SELinuxContext_abcAccept%3D0.target
[Unit]
Description=Foo
[Socket]
SELinuxContext=abc
Accept=0
