# tests/auto/SettingRequires/2.31/good/Socket_SELinuxContext_abcAccept%3Dtrue.target
[Unit]
Description=Foo
[Socket]
SELinuxContext=abc
Accept=true
