# systemdlint/tests/auto/OptionDeprecated/2.46/bad/TrafficControlQueueingDiscipline_NetworkEmulatorLossRate_0%25.target
[Unit]
Description=Foo
[TrafficControlQueueingDiscipline]
NetworkEmulatorLossRate=123

