# systemdlint/tests/auto/InvalidValue/2.48/good/DHCPv6PrefixDelegation_ManageTemporaryAddress_abc.target
[Unit]
Description=Foo
[DHCPv6PrefixDelegation]
ManageTemporaryAddress=abc

