# systemdlint/tests/auto/InvalidValue/2.48/good/DHCPv6_UseAddress_false.target
[Unit]
Description=Foo
[DHCPv6]
UseAddress=false

