# systemdlint/tests/auto/InvalidValue/2.48/good/IPv6AcceptRA_RouteAllowList_abc.target
[Unit]
Description=Foo
[IPv6AcceptRA]
RouteAllowList=abc

