# systemdlint/tests/auto/InvalidValue/2.48/good/IPv6AcceptRA_PrefixDenyList_abc.target
[Unit]
Description=Foo
[IPv6AcceptRA]
PrefixDenyList=abc

