# systemdlint/tests/auto/InvalidValue/2.48/good/Link_TransmitQueueLength_0.target
[Unit]
Description=Foo
[Link]
TransmitQueueLength=0

