# systemdlint/tests/auto/InvalidValue/2.48/good/Link_ActivationPolicy_bound.target
[Unit]
Description=Foo
[Link]
ActivationPolicy=bound

