# systemdlint/tests/auto/InvalidValue/2.48/good/Link_ReceiveQueues_4096.target
[Unit]
Description=Foo
[Link]
ReceiveQueues=4096

