# systemdlint/tests/auto/InvalidValue/2.48/good/Link_ActivationPolicy_manual.target
[Unit]
Description=Foo
[Link]
ActivationPolicy=manual

