# systemdlint/tests/auto/InvalidValue/2.48/good/Link_ReceiveQueues_1.target
[Unit]
Description=Foo
[Link]
ReceiveQueues=1

