# systemdlint/tests/auto/InvalidValue/2.48/good/Link_ActivationPolicy_always-down.target
[Unit]
Description=Foo
[Link]
ActivationPolicy=always-down

