# systemdlint/tests/auto/InvalidValue/2.48/good/RoutingPolicyRule_Type_prohibit.target
[Unit]
Description=Foo
[RoutingPolicyRule]
Type=prohibit

