# systemdlint/tests/auto/InvalidValue/2.48/good/RoutingPolicyRule_Family_both.target
[Unit]
Description=Foo
[RoutingPolicyRule]
Family=both

