# systemdlint/tests/auto/InvalidValue/2.48/good/Scope_ManagedOOMPreference_avoid.target
[Unit]
Description=Foo
[Scope]
ManagedOOMPreference=avoid

