# systemdlint/tests/auto/InvalidValue/2.48/good/Service_ManagedOOMPreference_omit.target
[Unit]
Description=Foo
[Service]
ManagedOOMPreference=omit

