# systemdlint/tests/auto/InvalidValue/2.48/good/Service_NoExecPaths_abc.target
[Unit]
Description=Foo
[Service]
NoExecPaths=abc

