# systemdlint/tests/auto/InvalidValue/2.48/good/Slice_ManagedOOMPreference_omit.target
[Unit]
Description=Foo
[Slice]
ManagedOOMPreference=omit

