# systemdlint/tests/auto/InvalidValue/2.48/good/Socket_ManagedOOMPreference_avoid.target
[Unit]
Description=Foo
[Socket]
ManagedOOMPreference=avoid

