# systemdlint/tests/auto/InvalidValue/2.48/good/Socket_ExecPaths_abc.target
[Unit]
Description=Foo
[Socket]
ExecPaths=abc

