# systemdlint/tests/auto/InvalidValue/2.48/good/Socket_ManagedOOMPreference_none.target
[Unit]
Description=Foo
[Socket]
ManagedOOMPreference=none

