# systemdlint/tests/auto/InvalidValue/2.48/good/Swap_ManagedOOMPreference_avoid.target
[Unit]
Description=Foo
[Swap]
ManagedOOMPreference=avoid

