# systemdlint/tests/auto/InvalidValue/2.48/good/Swap_ManagedOOMPreference_none.target
[Unit]
Description=Foo
[Swap]
ManagedOOMPreference=none

