# systemdlint/tests/auto/InvalidValue/2.48/good/Time_ConnectionRetrySec_9999999999999.target
[Unit]
Description=Foo
[Time]
ConnectionRetrySec=9999999999999

