# systemdlint/tests/auto/InvalidValue/2.48/good/Unit_ConditionCPUFeature_f16c.target
[Unit]
Description=Foo
Description=abc
ConditionCPUFeature=f16c

