# systemdlint/tests/auto/InvalidValue/2.48/good/Unit_ConditionCPUFeature_fma3.target
[Unit]
Description=Foo
Description=abc
ConditionCPUFeature=fma3

