# systemdlint/tests/auto/InvalidValue/2.48/good/Unit_ConditionCPUFeature_adx.target
[Unit]
Description=Foo
Description=abc
ConditionCPUFeature=adx

