# systemdlint/tests/auto/InvalidValue/2.48/good/Unit_ConditionCPUFeature_fxsr.target
[Unit]
Description=Foo
Description=abc
ConditionCPUFeature=fxsr

