# systemdlint/tests/auto/InvalidValue/2.48/good/Unit_ConditionCPUFeature_pge.target
[Unit]
Description=Foo
Description=abc
ConditionCPUFeature=pge

