# systemdlint/tests/auto/InvalidValue/2.48/good/Unit_ConditionCPUFeature_constant_tsc.target
[Unit]
Description=Foo
Description=abc
ConditionCPUFeature=constant_tsc

