# systemdlint/tests/auto/InvalidValue/2.35/good/Mount_IPAccounting_1.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
IPAccounting=1

