# systemdlint/tests/auto/InvalidValue/2.35/good/Mount_ConfigurationDirectoryMode_9999999999999.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
ConfigurationDirectoryMode=9999999999999

