# systemdlint/tests/auto/InvalidValue/2.35/good/Mount_IPAddressDeny_255.255.255.255--32.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
IPAddressDeny=255.255.255.255/32

