# systemdlint/tests/auto/InvalidValue/2.35/good/Service_StateDirectoryMode_0.target
[Unit]
Description=Foo
[Service]
StateDirectoryMode=0

