# systemdlint/tests/auto/InvalidValue/2.35/good/Service_CacheDirectoryMode_0.target
[Unit]
Description=Foo
[Service]
CacheDirectoryMode=0

