# systemdlint/tests/auto/InvalidValue/2.35/good/Service_LockPersonality_off.target
[Unit]
Description=Foo
[Service]
LockPersonality=off

