# systemdlint/tests/auto/InvalidValue/2.35/good/Service_KeyringMode_inherit.target
[Unit]
Description=Foo
[Service]
KeyringMode=inherit

