# systemdlint/tests/auto/InvalidValue/2.35/good/Service_StateDirectoryMode_9999999999999.target
[Unit]
Description=Foo
[Service]
StateDirectoryMode=9999999999999

