# systemdlint/tests/auto/InvalidValue/2.35/good/Service_RuntimeDirectoryMode_9999999999999.target
[Unit]
Description=Foo
[Service]
RuntimeDirectoryMode=9999999999999

