# systemdlint/tests/auto/InvalidValue/2.35/good/Service_UnsetEnvironment_abc.target
[Unit]
Description=Foo
[Service]
UnsetEnvironment=abc

