# systemdlint/tests/auto/InvalidValue/2.35/good/Service_LockPersonality_1.target
[Unit]
Description=Foo
[Service]
LockPersonality=1

