# systemdlint/tests/auto/InvalidValue/2.35/good/Service_RuntimeDirectoryPreserve_false.target
[Unit]
Description=Foo
[Service]
RuntimeDirectoryPreserve=false

