# systemdlint/tests/auto/InvalidValue/2.35/good/Service_IPAddressAllow_255.255.255.255.target
[Unit]
Description=Foo
[Service]
IPAddressAllow=255.255.255.255

