# systemdlint/tests/auto/InvalidValue/2.35/good/Service_RuntimeDirectoryPreserve_restart.target
[Unit]
Description=Foo
[Service]
RuntimeDirectoryPreserve=restart

