# systemdlint/tests/auto/InvalidValue/2.35/good/Service_RuntimeDirectoryPreserve_on.target
[Unit]
Description=Foo
[Service]
RuntimeDirectoryPreserve=on

