# systemdlint/tests/auto/InvalidValue/2.35/good/Service_IPAddressDeny_255.255.255.255.target
[Unit]
Description=Foo
[Service]
IPAddressDeny=255.255.255.255

