# systemdlint/tests/auto/InvalidValue/2.35/good/Service_IPAccounting_false.target
[Unit]
Description=Foo
[Service]
IPAccounting=false

