# systemdlint/tests/auto/InvalidValue/2.35/good/Service_IPAccounting_0.target
[Unit]
Description=Foo
[Service]
IPAccounting=0

