# systemdlint/tests/auto/InvalidValue/2.35/good/Service_RuntimeDirectoryPreserve_off.target
[Unit]
Description=Foo
[Service]
RuntimeDirectoryPreserve=off

