# systemdlint/tests/auto/InvalidValue/2.35/good/Service_RuntimeDirectoryMode_0.target
[Unit]
Description=Foo
[Service]
RuntimeDirectoryMode=0

