# systemdlint/tests/auto/InvalidValue/2.35/good/Service_RuntimeDirectory_abc.target
[Unit]
Description=Foo
[Service]
RuntimeDirectory=abc

