# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_StateDirectory_abc.target
[Unit]
Description=Foo
[Socket]
StateDirectory=abc

