# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_KeyringMode_shared.target
[Unit]
Description=Foo
[Socket]
KeyringMode=shared

