# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_StateDirectoryMode_0.target
[Unit]
Description=Foo
[Socket]
StateDirectoryMode=0

