# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_LockPersonality_true.target
[Unit]
Description=Foo
[Socket]
LockPersonality=true

