# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_RuntimeDirectoryPreserve_restart.target
[Unit]
Description=Foo
[Socket]
RuntimeDirectoryPreserve=restart

