# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_CacheDirectoryMode_0.target
[Unit]
Description=Foo
[Socket]
CacheDirectoryMode=0

