# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_LockPersonality_1.target
[Unit]
Description=Foo
[Socket]
LockPersonality=1

