# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_UnsetEnvironment_abc.target
[Unit]
Description=Foo
[Socket]
UnsetEnvironment=abc

