# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_IPAccounting_yes.target
[Unit]
Description=Foo
[Socket]
IPAccounting=yes

