# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_RuntimeDirectoryPreserve_on.target
[Unit]
Description=Foo
[Socket]
RuntimeDirectoryPreserve=on

