# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_CacheDirectoryMode_9999999999999.target
[Unit]
Description=Foo
[Socket]
CacheDirectoryMode=9999999999999

