# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_LockPersonality_off.target
[Unit]
Description=Foo
[Socket]
LockPersonality=off

