# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_LogsDirectory_abc.target
[Unit]
Description=Foo
[Socket]
LogsDirectory=abc

