# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_IPAccounting_1.target
[Unit]
Description=Foo
[Socket]
IPAccounting=1

