# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_IPAccounting_no.target
[Unit]
Description=Foo
[Socket]
IPAccounting=no

