# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_RuntimeDirectoryPreserve_1.target
[Unit]
Description=Foo
[Socket]
RuntimeDirectoryPreserve=1

