# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_RuntimeDirectoryPreserve_no.target
[Unit]
Description=Foo
[Socket]
RuntimeDirectoryPreserve=no

