# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_ConfigurationDirectoryMode_9999999999999.target
[Unit]
Description=Foo
[Socket]
ConfigurationDirectoryMode=9999999999999

