# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_ConfigurationDirectory_abc.target
[Unit]
Description=Foo
[Socket]
ConfigurationDirectory=abc

