# systemdlint/tests/auto/InvalidValue/2.35/good/Socket_CacheDirectory_abc.target
[Unit]
Description=Foo
[Socket]
CacheDirectory=abc

