# systemdlint/tests/auto/InvalidValue/2.35/good/Swap_IPAccounting_0.target
[Unit]
Description=Foo
[Swap]
IPAccounting=0

