# systemdlint/tests/auto/InvalidValue/2.35/good/Swap_IPAccounting_off.target
[Unit]
Description=Foo
[Swap]
IPAccounting=off

