# systemdlint/tests/auto/InvalidValue/2.35/good/Swap_RuntimeDirectoryPreserve_1.target
[Unit]
Description=Foo
[Swap]
RuntimeDirectoryPreserve=1

