# systemdlint/tests/auto/InvalidValue/2.35/good/Swap_RuntimeDirectoryMode_0.target
[Unit]
Description=Foo
[Swap]
RuntimeDirectoryMode=0

