# systemdlint/tests/auto/InvalidValue/2.35/good/Swap_LockPersonality_on.target
[Unit]
Description=Foo
[Swap]
LockPersonality=on

