# systemdlint/tests/auto/InvalidValue/2.35/good/Swap_CacheDirectoryMode_0.target
[Unit]
Description=Foo
[Swap]
CacheDirectoryMode=0

