# systemdlint/tests/auto/InvalidValue/2.35/good/Swap_CacheDirectory_abc.target
[Unit]
Description=Foo
[Swap]
CacheDirectory=abc

