# systemdlint/tests/auto/InvalidValue/2.35/good/Swap_LockPersonality_true.target
[Unit]
Description=Foo
[Swap]
LockPersonality=true

