# systemdlint/tests/auto/InvalidValue/2.35/good/Swap_IPAccounting_true.target
[Unit]
Description=Foo
[Swap]
IPAccounting=true

