# systemdlint/tests/auto/InvalidValue/2.35/good/Swap_RuntimeDirectoryPreserve_true.target
[Unit]
Description=Foo
[Swap]
RuntimeDirectoryPreserve=true

