# systemdlint/tests/auto/InvalidValue/2.35/good/Swap_StateDirectory_abc.target
[Unit]
Description=Foo
[Swap]
StateDirectory=abc

