# systemdlint/tests/auto/InvalidValue/2.35/good/Swap_KeyringMode_private.target
[Unit]
Description=Foo
[Swap]
KeyringMode=private

