# systemdlint/tests/auto/InvalidValue/2.35/good/Swap_KeyringMode_inherit.target
[Unit]
Description=Foo
[Swap]
KeyringMode=inherit

