# systemdlint/tests/auto/InvalidValue/2.35/good/Swap_RuntimeDirectoryPreserve_off.target
[Unit]
Description=Foo
[Swap]
RuntimeDirectoryPreserve=off

