# systemdlint/tests/auto/InvalidValue/2.35/good/Swap_LockPersonality_0.target
[Unit]
Description=Foo
[Swap]
LockPersonality=0

