# systemdlint/tests/auto/InvalidValue/2.35/good/Swap_RuntimeDirectoryPreserve_no.target
[Unit]
Description=Foo
[Swap]
RuntimeDirectoryPreserve=no

