# systemdlint/tests/auto/InvalidValue/2.40/good/Mount_MemoryMin_0M.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
MemoryMin=0M

