# systemdlint/tests/auto/InvalidValue/2.40/good/Mount_LogRateLimitBurst_19s.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
LogRateLimitBurst=19s

