# systemdlint/tests/auto/InvalidValue/2.40/good/Mount_MemoryMin_9999999999999T.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
MemoryMin=9999999999999T

