# systemdlint/tests/auto/InvalidValue/2.40/good/Mount_LogRateLimitBurst_22hr.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
LogRateLimitBurst=22hr

