# systemdlint/tests/auto/InvalidValue/2.40/good/Mount_LogRateLimitIntervalSec_26months.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
LogRateLimitIntervalSec=26months

