# systemdlint/tests/auto/InvalidValue/2.40/good/Mount_LogRateLimitIntervalSec_4ms.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
LogRateLimitIntervalSec=4ms

