# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_27msec.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=27msec

