# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_12month.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=12month

