# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_18sec.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=18sec

