# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_8weeks.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=8weeks

