# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_15min.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=15min

