# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_infinity.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=infinity

