# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_14M.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=14M

