# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_18year.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=18year

