# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_5hour.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=5hour

