# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_6h.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=6h

