# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_1%C2%B5s.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=1µs

