# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_28h.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=28h

