# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_24years.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=24years

