# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_0days.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=0days

