# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_3hours.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=3hours

