# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_28h.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=28h

