# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_22hr.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=22hr

