# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_10day.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=10day

