# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_25s.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=25s

