# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_26months.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=26months

