# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_9month.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=9month

