# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_infinity.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=infinity

