# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_22hr.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=22hr

