# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_19d.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=19d

