# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_17minutes.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=17minutes

