# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_20w.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=20w

