# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_26year.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=26year

