# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_21usec.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=21usec

