# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_11y.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=11y

