# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_28%C2%B5s.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=28µs

