# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_11y.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=11y

