# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_7seconds.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=7seconds

