# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_13minute.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=13minute

