# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_7seconds.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=7seconds

