# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_0sec.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=0sec

