# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_20usec.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=20usec

