# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_16seconds.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=16seconds

