# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_13minute.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=13minute

