# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_18year.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=18year

