# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_27y.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=27y

