# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_3hours.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=3hours

