# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_10M.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=10M

