# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_13months.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=13months

