# systemdlint/tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_6days.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=6days

