# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_13minute.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=13minute

