# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_infinity.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=infinity

