# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_6days.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=6days

