# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_26months.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=26months

