# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_7seconds.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=7seconds

