# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_27msec.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=27msec

