# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_10day.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=10day

