# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_10day.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=10day

