# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_15ms.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=15ms

