# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_MemoryMin_0B.target
[Unit]
Description=Foo
[Socket]
MemoryMin=0B

