# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_9month.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=9month

