# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_8weeks.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=8weeks

