# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_16week.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=16week

