# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_25s.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=25s

