# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_24years.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=24years

