# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_21usec.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=21usec

