# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_9month.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=9month

