# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_27y.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=27y

