# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_21us.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=21us

