# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_24years.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=24years

