# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_9min.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=9min

