# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_19s.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=19s

