# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_8minute.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=8minute

