# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_20w.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=20w

