# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_11y.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=11y

