# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_19s.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=19s

