# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_21usec.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=21usec

