# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_6days.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=6days

