# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_2d.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=2d

