# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_7minutes.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=7minutes

