# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_18year.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=18year

