# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_4hour.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=4hour

