# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_14M.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=14M

