# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_25years.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=25years

