# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_13months.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=13months

