# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_8minute.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=8minute

