# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_5hour.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=5hour

