# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_23m.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=23m

