# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_0sec.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=0sec

