# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_MemoryMin_9999999999999K.target
[Unit]
Description=Foo
[Socket]
MemoryMin=9999999999999K

