# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_16week.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=16week

