# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_3hours.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=3hours

