# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_12month.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=12month

