# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_5hour.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=5hour

