# systemdlint/tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_4ms.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=4ms

