# systemdlint/tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitIntervalSec_24years.target
[Unit]
Description=Foo
[Swap]
LogRateLimitIntervalSec=24years

