# systemdlint/tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitBurst_11y.target
[Unit]
Description=Foo
[Swap]
LogRateLimitBurst=11y

