# systemdlint/tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitIntervalSec_19d.target
[Unit]
Description=Foo
[Swap]
LogRateLimitIntervalSec=19d

