# systemdlint/tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitBurst_20usec.target
[Unit]
Description=Foo
[Swap]
LogRateLimitBurst=20usec

