# systemdlint/tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitIntervalSec_4ms.target
[Unit]
Description=Foo
[Swap]
LogRateLimitIntervalSec=4ms

