# systemdlint/tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitIntervalSec_9min.target
[Unit]
Description=Foo
[Swap]
LogRateLimitIntervalSec=9min

