# systemdlint/tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitBurst_25s.target
[Unit]
Description=Foo
[Swap]
LogRateLimitBurst=25s

