# systemdlint/tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitBurst_19s.target
[Unit]
Description=Foo
[Swap]
LogRateLimitBurst=19s

