# systemdlint/tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitIntervalSec_26months.target
[Unit]
Description=Foo
[Swap]
LogRateLimitIntervalSec=26months

