# systemdlint/tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitBurst_6h.target
[Unit]
Description=Foo
[Swap]
LogRateLimitBurst=6h

