# systemdlint/tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitIntervalSec_7minutes.target
[Unit]
Description=Foo
[Swap]
LogRateLimitIntervalSec=7minutes

