# systemdlint/tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitBurst_7minutes.target
[Unit]
Description=Foo
[Swap]
LogRateLimitBurst=7minutes

