# systemdlint/tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitIntervalSec_0sec.target
[Unit]
Description=Foo
[Swap]
LogRateLimitIntervalSec=0sec

