# systemdlint/tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitIntervalSec_10day.target
[Unit]
Description=Foo
[Swap]
LogRateLimitIntervalSec=10day

