# systemdlint/tests/auto/InvalidValue/2.40/good/Swap_MemoryMin_0B.target
[Unit]
Description=Foo
[Swap]
MemoryMin=0B

