# systemdlint/tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitIntervalSec_10M.target
[Unit]
Description=Foo
[Swap]
LogRateLimitIntervalSec=10M

