# systemdlint/tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitBurst_15ms.target
[Unit]
Description=Foo
[Swap]
LogRateLimitBurst=15ms

