# systemdlint/tests/auto/OptionTooNew/2.47/good/FlowQueuePIE_PacketLimit_0.target
[Unit]
Description=Foo
[FlowQueuePIE]
PacketLimit=123

