# systemdlint/tests/auto/OptionTooNew/2.47/good/Network_DHCPv6PrefixDelegation_yes.target
[Unit]
Description=Foo
[Network]
DHCPv6PrefixDelegation=123

