# systemdlint/tests/auto/OptionTooNew/2.47/good/OOM_DefaultMemoryPressureLimit_abc.target
[Unit]
Description=Foo
[OOM]
DefaultMemoryPressureLimit=123

