# systemdlint/tests/auto/OptionTooNew/2.47/good/Socket_ManagedOOMMemoryPressureLimitPercent_abc.target
[Unit]
Description=Foo
[Socket]
ManagedOOMMemoryPressureLimitPercent=123

