# systemdlint/tests/auto/OptionTooNew/2.47/good/Socket_SystemCallLog_abc.target
[Unit]
Description=Foo
[Socket]
SystemCallLog=123

