# systemdlint/tests/auto/OptionTooNew/2.47/good/Socket_RootImageOptions_%20root%20root-secondary%20home%20srv%20esp%20xbootldr%20tmp%20var%20usr.target
[Unit]
Description=Foo
[Socket]
RootImageOptions=123

