# systemdlint/tests/auto/OptionTooNew/2.47/good/Socket_MountImages_abc.target
[Unit]
Description=Foo
[Socket]
MountImages=123

