# systemdlint/tests/auto/OptionTooNew/2.47/good/Swap_SystemCallLog_abc.target
[Unit]
Description=Foo
[Swap]
SystemCallLog=123

