# systemdlint/tests/auto/OptionTooNew/2.47/good/VXLAN_Independent_yes.target
[Unit]
Description=Foo
[VXLAN]
Independent=123

