# tests/auto/OptionDeprecated/2.42/bad/Socket_BlockIOWriteBandwidthMax_--dev--foo%200K.target
[Unit]
Description=Foo
[Socket]
BlockIOWriteBandwidthMax=123

