# tests/auto/OptionDeprecated/2.42/bad/Swap_BlockStartupIOWeight_1.target
[Unit]
Description=Foo
[Swap]
BlockStartupIOWeight=123

