# systemdlint/tests/auto/InvalidValue/2.27/good/Manager_CrashChVT_abc.target
[Unit]
Description=Foo
[Manager]
CrashChVT=abc

