# systemdlint/tests/auto/InvalidValue/2.42/bad/Service_CPUQuotaPeriodSec_3decades.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=3decades

