# systemdlint/tests/auto/InvalidValue/2.42/bad/Service_NetworkNamespacePath_1.target
[Unit]
Description=Foo
[Service]
NetworkNamespacePath=1

