# systemdlint/tests/auto/InvalidValue/2.42/bad/Service_ProtectHostname_foo.target
[Unit]
Description=Foo
[Service]
ProtectHostname=foo

