# systemdlint/tests/auto/InvalidValue/2.42/bad/Service_CPUQuotaPeriodSec_True.target
[Unit]
Description=Foo
[Service]
CPUQuotaPeriodSec=True

