# systemdlint/tests/auto/InvalidValue/2.42/bad/Service_ProtectHostname_bar.target
[Unit]
Description=Foo
[Service]
ProtectHostname=bar

