# systemdlint/tests/auto/InvalidValue/2.42/bad/Service_NetworkNamespacePath_tmp.target
[Unit]
Description=Foo
[Service]
NetworkNamespacePath=tmp

