# systemdlint/tests/auto/InvalidValue/2.42/bad/Socket_CPUQuotaPeriodSec_True.target
[Unit]
Description=Foo
[Socket]
CPUQuotaPeriodSec=True

