# systemdlint/tests/auto/InvalidValue/2.42/bad/Socket_NetworkNamespacePath_..--..--.target
[Unit]
Description=Foo
[Socket]
NetworkNamespacePath=../../

