# systemdlint/tests/auto/InvalidValue/2.42/bad/Socket_ProtectHostname_baz.target
[Unit]
Description=Foo
[Socket]
ProtectHostname=baz

