# systemdlint/tests/auto/InvalidValue/2.42/bad/Socket_CPUQuotaPeriodSec_0.5d.target
[Unit]
Description=Foo
[Socket]
CPUQuotaPeriodSec=0.5d

