# systemdlint/tests/auto/InvalidValue/2.42/bad/Swap_CPUQuotaPeriodSec_1p.target
[Unit]
Description=Foo
[Swap]
CPUQuotaPeriodSec=1p

