# systemdlint/tests/auto/InvalidValue/2.42/bad/Swap_ProtectHostname_bar.target
[Unit]
Description=Foo
[Swap]
ProtectHostname=bar

