# systemdlint/tests/auto/InvalidValue/2.42/bad/Unit_ConditionCPUs_-1.target
[Unit]
Description=Foo
Description=abc
ConditionCPUs=-1

