# tests/manual/ConflictingOptions/2.31/good/Unit_ReloadPropagatedFrom.target
[Unit]
Description=Before and After
ReloadPropagatedFrom=foo.service