# tests/manual/ConflictingOptions/2.31/good/Unit_Requires.target
[Unit]
Description=Before and After
Requires=foo.service