# tests/manual/ConflictingOptions/2.31/good/Unit_After.target
[Unit]
Description=Before and After
After=foo.service