# tests/manual/ConflictingOptions/2.31/good/Unit_Wants.target
[Unit]
Description=Before and After
Wants=foo.service