# systemdlint/tests/auto/OptionTooNew/2.32/good/Socket_CPUWeight_1.target
[Unit]
Description=Foo
[Socket]
CPUWeight=123

