# systemdlint/tests/auto/OptionTooNew/2.32/good/Swap_ProtectControlGroups_yes.target
[Unit]
Description=Foo
[Swap]
ProtectControlGroups=123

