# systemdlint/tests/auto/InvalidValue/2.35/bad/Mount_CacheDirectoryMode_10000000000000.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
CacheDirectoryMode=10000000000000

