# systemdlint/tests/auto/InvalidValue/2.35/bad/Mount_IPAddressAllow_255.256.255.255.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
IPAddressAllow=255.256.255.255

