# systemdlint/tests/auto/InvalidValue/2.35/bad/Service_KeyringMode_1.target
[Unit]
Description=Foo
[Service]
KeyringMode=1

