# systemdlint/tests/auto/InvalidValue/2.35/bad/Service_CacheDirectoryMode_-1.target
[Unit]
Description=Foo
[Service]
CacheDirectoryMode=-1

