# systemdlint/tests/auto/InvalidValue/2.35/bad/Service_KeyringMode_True.target
[Unit]
Description=Foo
[Service]
KeyringMode=True

