# systemdlint/tests/auto/InvalidValue/2.35/bad/Service_IPAccounting_baz.target
[Unit]
Description=Foo
[Service]
IPAccounting=baz

