# systemdlint/tests/auto/InvalidValue/2.35/bad/Service_LockPersonality_foo.target
[Unit]
Description=Foo
[Service]
LockPersonality=foo

