# systemdlint/tests/auto/InvalidValue/2.35/bad/Service_IPAddressDeny_256.255.255.255.target
[Unit]
Description=Foo
[Service]
IPAddressDeny=256.255.255.255

