# systemdlint/tests/auto/InvalidValue/2.35/bad/Service_IPAddressDeny_0.0.0..target
[Unit]
Description=Foo
[Service]
IPAddressDeny=0.0.0.

