# systemdlint/tests/auto/InvalidValue/2.35/bad/Service_IPAccounting_bar.target
[Unit]
Description=Foo
[Service]
IPAccounting=bar

