# systemdlint/tests/auto/InvalidValue/2.35/bad/Service_IPAccounting_foo.target
[Unit]
Description=Foo
[Service]
IPAccounting=foo

