# systemdlint/tests/auto/InvalidValue/2.35/bad/Service_RuntimeDirectoryMode_10000000000000.target
[Unit]
Description=Foo
[Service]
RuntimeDirectoryMode=10000000000000

