# systemdlint/tests/auto/InvalidValue/2.35/bad/Service_LockPersonality_baz.target
[Unit]
Description=Foo
[Service]
LockPersonality=baz

