# systemdlint/tests/auto/InvalidValue/2.35/bad/Service_IPAddressDeny_abc.target
[Unit]
Description=Foo
[Service]
IPAddressDeny=abc

