# systemdlint/tests/auto/InvalidValue/2.35/bad/Service_IPAddressAllow_0--abc.target
[Unit]
Description=Foo
[Service]
IPAddressAllow=0/abc

