# systemdlint/tests/auto/InvalidValue/2.35/bad/Socket_StateDirectoryMode_10000000000000.target
[Unit]
Description=Foo
[Socket]
StateDirectoryMode=10000000000000

