# systemdlint/tests/auto/InvalidValue/2.35/bad/Socket_LockPersonality_foo.target
[Unit]
Description=Foo
[Socket]
LockPersonality=foo

