# systemdlint/tests/auto/InvalidValue/2.35/bad/Socket_IPAddressAllow_256.255.255.255.target
[Unit]
Description=Foo
[Socket]
IPAddressAllow=256.255.255.255

