# systemdlint/tests/auto/InvalidValue/2.35/bad/Socket_KeyringMode_True.target
[Unit]
Description=Foo
[Socket]
KeyringMode=True

