# systemdlint/tests/auto/InvalidValue/2.35/bad/Socket_IPAddressAllow_-1.1.1.1.target
[Unit]
Description=Foo
[Socket]
IPAddressAllow=-1.1.1.1

