# systemdlint/tests/auto/InvalidValue/2.35/bad/Socket_ConfigurationDirectoryMode_10000000000000.target
[Unit]
Description=Foo
[Socket]
ConfigurationDirectoryMode=10000000000000

