# systemdlint/tests/auto/InvalidValue/2.35/bad/Socket_KeyringMode_1.target
[Unit]
Description=Foo
[Socket]
KeyringMode=1

