# systemdlint/tests/auto/InvalidValue/2.35/bad/Socket_StateDirectoryMode_-1.target
[Unit]
Description=Foo
[Socket]
StateDirectoryMode=-1

