# systemdlint/tests/auto/InvalidValue/2.35/bad/Socket_IPAddressDeny_255.255.255.256.target
[Unit]
Description=Foo
[Socket]
IPAddressDeny=255.255.255.256

