# systemdlint/tests/auto/InvalidValue/2.35/bad/Socket_RuntimeDirectoryMode_-1.target
[Unit]
Description=Foo
[Socket]
RuntimeDirectoryMode=-1

