# systemdlint/tests/auto/InvalidValue/2.35/bad/Swap_IPAddressDeny_256.255.255.255.target
[Unit]
Description=Foo
[Swap]
IPAddressDeny=256.255.255.255

