# systemdlint/tests/auto/InvalidValue/2.35/bad/Swap_IPAddressAllow_256.255.255.255.target
[Unit]
Description=Foo
[Swap]
IPAddressAllow=256.255.255.255

