# systemdlint/tests/auto/InvalidValue/2.35/bad/Swap_IPAddressDeny_0.0.0..target
[Unit]
Description=Foo
[Swap]
IPAddressDeny=0.0.0.

