# systemdlint/tests/auto/InvalidValue/2.35/bad/Swap_IPAddressDeny_-1.1.1.1.target
[Unit]
Description=Foo
[Swap]
IPAddressDeny=-1.1.1.1

