# systemdlint/tests/auto/InvalidValue/2.35/bad/Swap_KeyringMode_True.target
[Unit]
Description=Foo
[Swap]
KeyringMode=True

