# systemdlint/tests/auto/InvalidValue/2.35/bad/Swap_KeyringMode_1.target
[Unit]
Description=Foo
[Swap]
KeyringMode=1

