# systemdlint/tests/auto/InvalidValue/2.35/bad/Swap_IPAddressAllow_0.0.0..target
[Unit]
Description=Foo
[Swap]
IPAddressAllow=0.0.0.

