# systemdlint/tests/auto/InvalidValue/2.35/bad/Swap_IPAddressDeny_255.255.255.255--33.target
[Unit]
Description=Foo
[Swap]
IPAddressDeny=255.255.255.255/33

