# systemdlint/tests/auto/InvalidValue/2.35/bad/Swap_ConfigurationDirectoryMode_-1.target
[Unit]
Description=Foo
[Swap]
ConfigurationDirectoryMode=-1

