# systemdlint/tests/auto/InvalidValue/2.35/bad/Swap_RuntimeDirectoryMode_-1.target
[Unit]
Description=Foo
[Swap]
RuntimeDirectoryMode=-1

