# systemdlint/tests/auto/InvalidValue/2.35/bad/Swap_IPAccounting_bar.target
[Unit]
Description=Foo
[Swap]
IPAccounting=bar

