# systemdlint/tests/auto/InvalidValue/2.35/bad/Swap_ConfigurationDirectoryMode_10000000000000.target
[Unit]
Description=Foo
[Swap]
ConfigurationDirectoryMode=10000000000000

