# tests/auto/InvalidValue/2.44/bad/DHCPv6_PrefixDelegationHint_255.256.255.255.target
[Unit]
Description=Foo
[DHCPv6]
PrefixDelegationHint=255.256.255.255

