# tests/auto/InvalidValue/2.44/bad/DHCPv6_PrefixDelegationHint_0--abc.target
[Unit]
Description=Foo
[DHCPv6]
PrefixDelegationHint=0/abc

