# tests/auto/InvalidValue/2.44/bad/DHCPServer_EmitSIP_bar.target
[Unit]
Description=Foo
[DHCPServer]
EmitSIP=bar

