# tests/auto/InvalidValue/2.44/bad/LifetimeSec_Route_10000000000000.target
[Unit]
Description=Foo
[LifetimeSec]
Route=10000000000000

