# tests/auto/InvalidValue/2.44/bad/Scope_RuntimeMaxSec_0.5d.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=0.5d

