# tests/auto/InvalidValue/2.44/bad/Scope_RuntimeMaxSec_2potatoes.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=2potatoes

