# tests/auto/InvalidValue/2.44/bad/Scope_RuntimeMaxSec_True.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=True

