# tests/auto/InvalidValue/2.44/bad/Service_ProtectKernelLogs_bar.target
[Unit]
Description=Foo
[Service]
ProtectKernelLogs=bar

