# tests/auto/InvalidValue/2.44/bad/Service_ProtectKernelLogs_foo.target
[Unit]
Description=Foo
[Service]
ProtectKernelLogs=foo

