# tests/auto/InvalidValue/2.44/bad/Socket_ProtectKernelLogs_foo.target
[Unit]
Description=Foo
[Socket]
ProtectKernelLogs=foo

