# tests/auto/InvalidValue/2.44/bad/Swap_ProtectKernelLogs_foo.target
[Unit]
Description=Foo
[Swap]
ProtectKernelLogs=foo

