# tests/auto/InvalidValue/2.44/bad/TrafficControlQueueingDiscipline_NetworkEmulatorLossRate_-1%25.target
[Unit]
Description=Foo
[TrafficControlQueueingDiscipline]
NetworkEmulatorLossRate=-1%

