# tests/auto/InvalidValue/2.44/bad/TrafficControlQueueingDiscipline_NetworkEmulatorDelaySec_10000000000000.target
[Unit]
Description=Foo
[TrafficControlQueueingDiscipline]
NetworkEmulatorDelaySec=10000000000000

