# tests/auto/InvalidValue/2.44/bad/TrafficControlQueueingDiscipline_NetworkEmulatorLossRate_101%25.target
[Unit]
Description=Foo
[TrafficControlQueueingDiscipline]
NetworkEmulatorLossRate=101%

