# tests/workspace/ReferencedUnitNotFound/2.31/good/Install_RequiredBy_foo.service.target/etc/systemd/system/bar.service
[Unit]
Description=foo
After=test.mount
Before=abc.swap
PartOf=grub-common.service