# tests/auto/OptionTooNew/2.11/bad/Timer_Persistent_yes.target
[Unit]
Description=Foo
[Timer]
Persistent=123

