# systemdlint/tests/auto/InvalidValue/2.44/bad/DHCPServer_EmitSIP_baz.target
[Unit]
Description=Foo
[DHCPServer]
EmitSIP=baz

