# systemdlint/tests/auto/InvalidValue/2.44/bad/DHCPServer_EmitSIP_foo.target
[Unit]
Description=Foo
[DHCPServer]
EmitSIP=foo

