# systemdlint/tests/auto/InvalidValue/2.44/bad/Link_RxBufferSize_-1.target
[Unit]
Description=Foo
[Link]
RxBufferSize=-1

