# systemdlint/tests/auto/InvalidValue/2.44/bad/Link_RxBufferSize_10000000000000.target
[Unit]
Description=Foo
[Link]
RxBufferSize=10000000000000

