# systemdlint/tests/auto/InvalidValue/2.44/bad/Link_TxBufferSize_-1.target
[Unit]
Description=Foo
[Link]
TxBufferSize=-1

