# systemdlint/tests/auto/InvalidValue/2.44/bad/Link_TxBufferSize_10000000000000.target
[Unit]
Description=Foo
[Link]
TxBufferSize=10000000000000

