# systemdlint/tests/auto/InvalidValue/2.44/bad/NextHop_Gateway_256.255.255.255.target
[Unit]
Description=Foo
[NextHop]
Gateway=1.1.1.1/0
Gateway=256.255.255.255

