# systemdlint/tests/auto/InvalidValue/2.44/bad/Scope_RuntimeMaxSec_1p.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=1p

