# systemdlint/tests/auto/InvalidValue/2.44/bad/Scope_RuntimeMaxSec_3decades.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=3decades

