# systemdlint/tests/auto/InvalidValue/2.44/bad/Service_ProtectKernelLogs_baz.target
[Unit]
Description=Foo
[Service]
ProtectKernelLogs=baz

