# systemdlint/tests/auto/InvalidValue/2.44/bad/Socket_ProtectKernelLogs_bar.target
[Unit]
Description=Foo
[Socket]
ProtectKernelLogs=bar

