# systemdlint/tests/auto/InvalidValue/2.44/bad/Socket_ProtectKernelLogs_baz.target
[Unit]
Description=Foo
[Socket]
ProtectKernelLogs=baz

