# systemdlint/tests/auto/InvalidValue/2.44/bad/Swap_ProtectKernelLogs_baz.target
[Unit]
Description=Foo
[Swap]
ProtectKernelLogs=baz

