# systemdlint/tests/auto/InvalidValue/2.44/bad/Swap_ProtectKernelLogs_bar.target
[Unit]
Description=Foo
[Swap]
ProtectKernelLogs=bar

