# systemdlint/tests/auto/InvalidValue/2.44/bad/TrafficControlQueueingDiscipline_NetworkEmulatorDelayJitterSec_-1.target
[Unit]
Description=Foo
[TrafficControlQueueingDiscipline]
NetworkEmulatorDelayJitterSec=-1

