# systemdlint/tests/auto/InvalidValue/2.44/bad/TrafficControlQueueingDiscipline_NetworkEmulatorPacketLimit_-1.target
[Unit]
Description=Foo
[TrafficControlQueueingDiscipline]
NetworkEmulatorPacketLimit=-1

