# systemdlint/tests/auto/OptionTooNew/2.50/good/Automount_ExtraOptions_abc.target
[Unit]
Description=Foo
[Automount]
ExtraOptions=123

