# systemdlint/tests/auto/OptionTooNew/2.50/good/CAKE_FirewallMark_1.target
[Unit]
Description=Foo
[CAKE]
FirewallMark=123

