# systemdlint/tests/auto/OptionTooNew/2.50/good/CAKE_PriorityQueueingPreset_besteffort.target
[Unit]
Description=Foo
[CAKE]
PriorityQueueingPreset=123

