# systemdlint/tests/auto/OptionTooNew/2.50/good/DHCPServer_Router_abc.target
[Unit]
Description=Foo
[DHCPServer]
Router=123

