# systemdlint/tests/auto/OptionTooNew/2.50/good/Link_WakeOnLanPassword_abc.target
[Unit]
Description=Foo
[Link]
WakeOnLanPassword=123

