# systemdlint/tests/auto/OptionTooNew/2.50/good/Mount_StartupAllowedMemoryNodes_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
StartupAllowedMemoryNodes=123

