# systemdlint/tests/auto/OptionTooNew/2.50/good/Path_TriggerLimitBurst_0.target
[Unit]
Description=Foo
[Path]
TriggerLimitBurst=123

