# systemdlint/tests/auto/OptionTooNew/2.50/good/Service_ExecSearchPath_abc.target
[Unit]
Description=Foo
[Service]
ExecSearchPath=123

