# systemdlint/tests/auto/OptionTooNew/2.50/good/Service_AllowedMemoryNodes_abc.target
[Unit]
Description=Foo
[Service]
AllowedMemoryNodes=123

