# systemdlint/tests/auto/OptionTooNew/2.50/good/Service_RestrictNetworkInterfaces_abc.target
[Unit]
Description=Foo
[Service]
RestrictNetworkInterfaces=123

