# systemdlint/tests/auto/OptionTooNew/2.50/good/Service_ExitType_main.target
[Unit]
Description=Foo
[Service]
ExitType=123

