# systemdlint/tests/auto/OptionTooNew/2.50/good/Slice_StartupAllowedCPUs_abc.target
[Unit]
Description=Foo
[Slice]
StartupAllowedCPUs=123

