# systemdlint/tests/auto/OptionTooNew/2.50/good/Socket_RestrictNetworkInterfaces_abc.target
[Unit]
Description=Foo
[Socket]
RestrictNetworkInterfaces=123

