# systemdlint/tests/auto/OptionTooNew/2.50/good/Swap_TTYColumns_0.target
[Unit]
Description=Foo
[Swap]
TTYColumns=123

