# systemdlint/tests/auto/OptionTooNew/2.50/good/Unit_AssertCPUPressure_abc.target
[Unit]
Description=Foo
Description=abc
AssertCPUPressure=123

