# systemdlint/tests/auto/OptionTooNew/2.50/good/Unit_ConditionMemoryPressure_abc.target
[Unit]
Description=Foo
Description=abc
ConditionMemoryPressure=123

