# systemdlint/tests/auto/OptionTooNew/2.50/good/WireGuardPeer_AllowedIPs_abc.target
[Unit]
Description=Foo
[WireGuardPeer]
AllowedIPs=123

