# systemdlint/tests/auto/InvalidValue/2.31/bad/Address_PreferredLifetime_1.target
[Unit]
Description=Foo
[Address]
Address=abc
PreferredLifetime=1

