# systemdlint/tests/auto/InvalidValue/2.31/bad/Address_PreferredLifetime_True.target
[Unit]
Description=Foo
[Address]
Address=abc
PreferredLifetime=True

