# systemdlint/tests/auto/InvalidValue/2.31/bad/Address_Peer_255.256.255.255.target
[Unit]
Description=Foo
[Address]
Address=abc
Peer=255.256.255.255

