# systemdlint/tests/auto/InvalidValue/2.31/bad/Automount_DirectoryMode_True.target
[Unit]
Description=Foo
[Automount]
DirectoryMode=True

