# systemdlint/tests/auto/InvalidValue/2.31/bad/Automount_Where_..--..--.target
[Unit]
Description=Foo
[Automount]
Where=../../

