# systemdlint/tests/auto/InvalidValue/2.31/bad/Automount_TimeoutIdleSec_True.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=True

