# systemdlint/tests/auto/InvalidValue/2.31/bad/Automount_TimeoutIdleSec_1p.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=1p

