# systemdlint/tests/auto/InvalidValue/2.31/bad/Automount_TimeoutIdleSec_3decades.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=3decades

