# systemdlint/tests/auto/InvalidValue/2.31/bad/Automount_DirectoryMode_-1.target
[Unit]
Description=Foo
[Automount]
DirectoryMode=-1

