# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_TransmitHashPolicy_True.target
[Unit]
Description=Foo
[Bond]
TransmitHashPolicy=True

