# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_ARPAllTargets_1.target
[Unit]
Description=Foo
[Bond]
ARPAllTargets=1

