# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_PrimaryReselectPolicy_1.target
[Unit]
Description=Foo
[Bond]
PrimaryReselectPolicy=1

