# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_MIIMonitorSec_-1.target
[Unit]
Description=Foo
[Bond]
MIIMonitorSec=-1

