# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_AllSlavesActive_baz.target
[Unit]
Description=Foo
[Bond]
AllSlavesActive=baz

