# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_ResendIGMP_256.target
[Unit]
Description=Foo
[Bond]
ResendIGMP=256

