# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_Mode_True.target
[Unit]
Description=Foo
[Bond]
Mode=True

