# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_AllSlavesActive_foo.target
[Unit]
Description=Foo
[Bond]
AllSlavesActive=foo

