# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_ARPAllTargets_True.target
[Unit]
Description=Foo
[Bond]
ARPAllTargets=True

