# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_LACPTransmitRate_True.target
[Unit]
Description=Foo
[Bond]
LACPTransmitRate=True

