# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_AdSelect_True.target
[Unit]
Description=Foo
[Bond]
AdSelect=True

