# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_AdSelect_1.target
[Unit]
Description=Foo
[Bond]
AdSelect=1

