# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_FailOverMACPolicy_1.target
[Unit]
Description=Foo
[Bond]
FailOverMACPolicy=1

