# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_TransmitHashPolicy_1.target
[Unit]
Description=Foo
[Bond]
TransmitHashPolicy=1

