# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_MinLinks_10000000000000.target
[Unit]
Description=Foo
[Bond]
MinLinks=10000000000000

