# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_FailOverMACPolicy_True.target
[Unit]
Description=Foo
[Bond]
FailOverMACPolicy=True

