# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_PrimaryReselectPolicy_True.target
[Unit]
Description=Foo
[Bond]
PrimaryReselectPolicy=True

