# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_PacketsPerSlave_-1.target
[Unit]
Description=Foo
[Bond]
PacketsPerSlave=-1

