# systemdlint/tests/auto/InvalidValue/2.31/bad/Bond_ARPValidate_True.target
[Unit]
Description=Foo
[Bond]
ARPValidate=True

