# systemdlint/tests/auto/InvalidValue/2.31/bad/Bridge_Learning_baz.target
[Unit]
Description=Foo
[Bridge]
Learning=baz

