# systemdlint/tests/auto/InvalidValue/2.31/bad/Bridge_UseBPDU_bar.target
[Unit]
Description=Foo
[Bridge]
UseBPDU=bar

