# systemdlint/tests/auto/InvalidValue/2.31/bad/Bridge_MulticastToUnicast_bar.target
[Unit]
Description=Foo
[Bridge]
MulticastToUnicast=bar

