# systemdlint/tests/auto/InvalidValue/2.31/bad/Bridge_UnicastFlood_bar.target
[Unit]
Description=Foo
[Bridge]
UnicastFlood=bar

