# systemdlint/tests/auto/InvalidValue/2.31/bad/Bridge_NeighborSuppression_baz.target
[Unit]
Description=Foo
[Bridge]
NeighborSuppression=baz

