# systemdlint/tests/auto/InvalidValue/2.31/bad/Bridge_FastLeave_foo.target
[Unit]
Description=Foo
[Bridge]
FastLeave=foo

