# systemdlint/tests/auto/InvalidValue/2.31/bad/Bridge_HairPin_bar.target
[Unit]
Description=Foo
[Bridge]
HairPin=bar

