# systemdlint/tests/auto/InvalidValue/2.31/bad/Bridge_HairPin_foo.target
[Unit]
Description=Foo
[Bridge]
HairPin=foo

