# systemdlint/tests/auto/InvalidValue/2.31/bad/Bridge_MulticastToUnicast_foo.target
[Unit]
Description=Foo
[Bridge]
MulticastToUnicast=foo

