# systemdlint/tests/auto/InvalidValue/2.31/bad/Bridge_UseBPDU_baz.target
[Unit]
Description=Foo
[Bridge]
UseBPDU=baz

