# systemdlint/tests/auto/InvalidValue/2.31/bad/Bridge_Priority_64.target
[Unit]
Description=Foo
[Bridge]
Priority=64

