# systemdlint/tests/auto/InvalidValue/2.31/bad/Bridge_FastLeave_bar.target
[Unit]
Description=Foo
[Bridge]
FastLeave=bar

