# systemdlint/tests/auto/InvalidValue/2.31/bad/Bridge_FastLeave_baz.target
[Unit]
Description=Foo
[Bridge]
FastLeave=baz

