# systemdlint/tests/auto/InvalidValue/2.31/bad/Bridge_AllowPortToBeRoot_bar.target
[Unit]
Description=Foo
[Bridge]
AllowPortToBeRoot=bar

