# systemdlint/tests/auto/InvalidValue/2.31/bad/CAN_TripleSampling_baz.target
[Unit]
Description=Foo
[CAN]
TripleSampling=baz

