# systemdlint/tests/auto/InvalidValue/2.31/bad/CAN_TripleSampling_bar.target
[Unit]
Description=Foo
[CAN]
TripleSampling=bar

