# systemdlint/tests/auto/InvalidValue/2.31/bad/CAN_RestartSec_True.target
[Unit]
Description=Foo
[CAN]
RestartSec=True

