# systemdlint/tests/auto/InvalidValue/2.31/bad/CAN_Bitrate_10000000000000K.target
[Unit]
Description=Foo
[CAN]
Bitrate=10000000000000K

