# systemdlint/tests/auto/InvalidValue/2.31/bad/CAN_RestartSec_2potatoes.target
[Unit]
Description=Foo
[CAN]
RestartSec=2potatoes

