# systemdlint/tests/auto/InvalidValue/2.31/bad/CAN_Bitrate_5001.target
[Unit]
Description=Foo
[CAN]
Bitrate=5001

