# systemdlint/tests/auto/InvalidValue/2.31/bad/CAN_RestartSec_3decades.target
[Unit]
Description=Foo
[CAN]
RestartSec=3decades

