# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_UseRoutes_bar.target
[Unit]
Description=Foo
[DHCP]
UseRoutes=bar

