# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_ClientIdentifier_1.target
[Unit]
Description=Foo
[DHCP]
ClientIdentifier=1

