# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_UseHostname_baz.target
[Unit]
Description=Foo
[DHCP]
UseHostname=baz

