# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_ClientIdentifier_True.target
[Unit]
Description=Foo
[DHCP]
ClientIdentifier=True

