# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_UseDNS_bar.target
[Unit]
Description=Foo
[DHCP]
UseDNS=bar

