# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_UseRoutes_baz.target
[Unit]
Description=Foo
[DHCP]
UseRoutes=baz

