# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_SendHostname_foo.target
[Unit]
Description=Foo
[DHCP]
SendHostname=foo

