# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_UseMTU_foo.target
[Unit]
Description=Foo
[DHCP]
UseMTU=foo

