# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCPServer_DefaultLeaseTimeSec_1p.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=1p

