# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCPServer_EmitNTP_baz.target
[Unit]
Description=Foo
[DHCPServer]
EmitNTP=baz

