# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCPServer_PoolSize_-1.target
[Unit]
Description=Foo
[DHCPServer]
PoolSize=-1

