# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_RapidCommit_baz.target
[Unit]
Description=Foo
[DHCP]
RapidCommit=baz

