# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_Anonymize_bar.target
[Unit]
Description=Foo
[DHCP]
Anonymize=bar

