# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_UseNTP_baz.target
[Unit]
Description=Foo
[DHCP]
UseNTP=baz

