# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_UseHostname_bar.target
[Unit]
Description=Foo
[DHCP]
UseHostname=bar

