# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_UseNTP_bar.target
[Unit]
Description=Foo
[DHCP]
UseNTP=bar

