# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCPServer_EmitDNS_baz.target
[Unit]
Description=Foo
[DHCPServer]
EmitDNS=baz

