# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_RapidCommit_foo.target
[Unit]
Description=Foo
[DHCP]
RapidCommit=foo

