# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCPServer_DefaultLeaseTimeSec_0.5d.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=0.5d

