# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_UseDNS_baz.target
[Unit]
Description=Foo
[DHCP]
UseDNS=baz

