# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_SendHostname_baz.target
[Unit]
Description=Foo
[DHCP]
SendHostname=baz

