# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_UseHostname_foo.target
[Unit]
Description=Foo
[DHCP]
UseHostname=foo

