# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_CriticalConnection_foo.target
[Unit]
Description=Foo
[DHCP]
CriticalConnection=foo

