# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCPServer_DefaultLeaseTimeSec_True.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=True

