# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCPServer_EmitRouter_foo.target
[Unit]
Description=Foo
[DHCPServer]
EmitRouter=foo

