# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_UseDNS_foo.target
[Unit]
Description=Foo
[DHCP]
UseDNS=foo

