# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCPServer_EmitTimezone_bar.target
[Unit]
Description=Foo
[DHCPServer]
EmitTimezone=bar

