# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_UseMTU_baz.target
[Unit]
Description=Foo
[DHCP]
UseMTU=baz

