# systemdlint/tests/auto/InvalidValue/2.31/bad/DHCP_Anonymize_foo.target
[Unit]
Description=Foo
[DHCP]
Anonymize=foo

