# systemdlint/tests/auto/InvalidValue/2.31/bad/FooOverUDP_Port_0.target
[Unit]
Description=Foo
[FooOverUDP]
Port=0

