# systemdlint/tests/auto/InvalidValue/2.31/bad/FooOverUDP_Port_65535.target
[Unit]
Description=Foo
[FooOverUDP]
Port=65535

