# systemdlint/tests/auto/InvalidValue/2.31/bad/FooOverUDP_Protocol_256.target
[Unit]
Description=Foo
[FooOverUDP]
Protocol=256

