# systemdlint/tests/auto/InvalidValue/2.31/bad/GENEVE_DestinationPort_0.target
[Unit]
Description=Foo
[GENEVE]
DestinationPort=0

