# systemdlint/tests/auto/InvalidValue/2.31/bad/GENEVE_Remote_256.255.255.255.target
[Unit]
Description=Foo
[GENEVE]
Remote=256.255.255.255

