# systemdlint/tests/auto/InvalidValue/2.31/bad/GENEVE_UDPChecksum_foo.target
[Unit]
Description=Foo
[GENEVE]
UDPChecksum=foo

