# systemdlint/tests/auto/InvalidValue/2.31/bad/GENEVE_UDP6ZeroChecksumTx_bar.target
[Unit]
Description=Foo
[GENEVE]
UDP6ZeroChecksumTx=bar

