# systemdlint/tests/auto/InvalidValue/2.31/bad/GENEVE_TTL_256.target
[Unit]
Description=Foo
[GENEVE]
TTL=256

