# systemdlint/tests/auto/InvalidValue/2.31/bad/IPv6Prefix_PreferredLifetimeSec_True.target
[Unit]
Description=Foo
[IPv6Prefix]
PreferredLifetimeSec=True

