# systemdlint/tests/auto/InvalidValue/2.31/bad/IPv6PrefixDelegation_RouterLifetimeSec_True.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
RouterLifetimeSec=True

