# systemdlint/tests/auto/InvalidValue/2.31/bad/IPv6Prefix_ValidLifetimeSec_0.5d.target
[Unit]
Description=Foo
[IPv6Prefix]
ValidLifetimeSec=0.5d

