# systemdlint/tests/auto/InvalidValue/2.31/bad/IPv6PrefixDelegation_RouterPreference_1.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
RouterPreference=1

