# systemdlint/tests/auto/InvalidValue/2.31/bad/IPv6Prefix_PreferredLifetimeSec_0.5d.target
[Unit]
Description=Foo
[IPv6Prefix]
PreferredLifetimeSec=0.5d

