# systemdlint/tests/auto/InvalidValue/2.31/bad/IPVLAN_Flags_True.target
[Unit]
Description=Foo
[IPVLAN]
Flags=True

