# systemdlint/tests/auto/InvalidValue/2.31/bad/IPVLAN_Flags_1.target
[Unit]
Description=Foo
[IPVLAN]
Flags=1

