# systemdlint/tests/auto/InvalidValue/2.31/bad/IPVLAN_Mode_True.target
[Unit]
Description=Foo
[IPVLAN]
Mode=True

