# systemdlint/tests/auto/InvalidValue/2.31/bad/IPVLAN_Mode_1.target
[Unit]
Description=Foo
[IPVLAN]
Mode=1

