# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_Compress_baz.target
[Unit]
Description=Foo
[Journal]
Compress=baz

