# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_RateLimitIntervalSec_-1s.target
[Unit]
Description=Foo
[Journal]
RateLimitIntervalSec=-1s

