# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_SystemMaxFileSize_-1P.target
[Unit]
Description=Foo
[Journal]
SystemMaxFileSize=-1P

