# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_SyncIntervalSec_0.5d.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=0.5d

