# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_ForwardToWall_bar.target
[Unit]
Description=Foo
[Journal]
ForwardToWall=bar

