# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_ReadKMsg_foo.target
[Unit]
Description=Foo
[Journal]
ReadKMsg=foo

