# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_ForwardToConsole_foo.target
[Unit]
Description=Foo
[Journal]
ForwardToConsole=foo

