# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_Compress_foo.target
[Unit]
Description=Foo
[Journal]
Compress=foo

