# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_RuntimeMaxFileSize_-1T.target
[Unit]
Description=Foo
[Journal]
RuntimeMaxFileSize=-1T

