# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_LineMax_10000000000000T.target
[Unit]
Description=Foo
[Journal]
LineMax=10000000000000T

