# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_TTYPath_--tmp.target
[Unit]
Description=Foo
[Journal]
TTYPath=/tmp

