# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_ForwardToSyslog_foo.target
[Unit]
Description=Foo
[Journal]
ForwardToSyslog=foo

