# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_MaxRetentionSec_-1day.target
[Unit]
Description=Foo
[Journal]
MaxRetentionSec=-1day

