# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_RuntimeMaxUse_-1P.target
[Unit]
Description=Foo
[Journal]
RuntimeMaxUse=-1P

