# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_TTYPath_--devvvv.target
[Unit]
Description=Foo
[Journal]
TTYPath=/devvvv

