# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_ReadKMsg_bar.target
[Unit]
Description=Foo
[Journal]
ReadKMsg=bar

