# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_Compress_bar.target
[Unit]
Description=Foo
[Journal]
Compress=bar

