# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_RuntimeKeepFree_-1P.target
[Unit]
Description=Foo
[Journal]
RuntimeKeepFree=-1P

