# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_ForwardToConsole_baz.target
[Unit]
Description=Foo
[Journal]
ForwardToConsole=baz

