# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_SystemMaxFiles_-1G.target
[Unit]
Description=Foo
[Journal]
SystemMaxFiles=-1G

