# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_SyncIntervalSec_True.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=True

