# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_SystemMaxUse_-1T.target
[Unit]
Description=Foo
[Journal]
SystemMaxUse=-1T

