# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_MaxFileSec_10000000000000year.target
[Unit]
Description=Foo
[Journal]
MaxFileSec=10000000000000year

