# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_RuntimeMaxUse_10000000000000G.target
[Unit]
Description=Foo
[Journal]
RuntimeMaxUse=10000000000000G

