# systemdlint/tests/auto/InvalidValue/2.31/bad/Journal_SystemKeepFree_10000000000000K.target
[Unit]
Description=Foo
[Journal]
SystemKeepFree=10000000000000K

