# systemdlint/tests/auto/InvalidValue/2.31/bad/L2TP_PeerTunnelId_0.target
[Unit]
Description=Foo
[L2TP]
PeerTunnelId=0

