# systemdlint/tests/auto/InvalidValue/2.31/bad/L2TP_UDP6ZeroChecksumRx_bar.target
[Unit]
Description=Foo
[L2TP]
UDP6ZeroChecksumRx=bar

