# systemdlint/tests/auto/InvalidValue/2.31/bad/L2TPSession_Layer2SpecificHeader_True.target
[Unit]
Description=Foo
[L2TPSession]
Name=abc
Layer2SpecificHeader=True

