# systemdlint/tests/auto/InvalidValue/2.31/bad/L2TP_DestinationPort_65535.target
[Unit]
Description=Foo
[L2TP]
DestinationPort=65535

