# systemdlint/tests/auto/InvalidValue/2.31/bad/L2TP_Remote_abc.target
[Unit]
Description=Foo
[L2TP]
Remote=abc

