# systemdlint/tests/auto/InvalidValue/2.31/bad/L2TPSession_PeerSessionId_4294967296.target
[Unit]
Description=Foo
[L2TPSession]
Name=abc
PeerSessionId=4294967296

