# systemdlint/tests/auto/InvalidValue/2.31/bad/Link_RequiredForOnline_baz.target
[Unit]
Description=Foo
[Link]
RequiredForOnline=baz

