# systemdlint/tests/auto/InvalidValue/2.31/bad/Link_Multicast_bar.target
[Unit]
Description=Foo
[Link]
Multicast=bar

