# systemdlint/tests/auto/InvalidValue/2.31/bad/Link_RequiredForOnline_bar.target
[Unit]
Description=Foo
[Link]
RequiredForOnline=bar

