# systemdlint/tests/auto/InvalidValue/2.31/bad/Link_MTUBytes_3073.target
[Unit]
Description=Foo
[Link]
MACAddress=abc
MTUBytes=3073

