# systemdlint/tests/auto/InvalidValue/2.31/bad/Link_Unmanaged_baz.target
[Unit]
Description=Foo
[Link]
Unmanaged=baz

