# systemdlint/tests/auto/InvalidValue/2.31/bad/Link_Unmanaged_foo.target
[Unit]
Description=Foo
[Link]
Unmanaged=foo

