# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_IdleAction_1.target
[Unit]
Description=Foo
[Login]
IdleAction=1

