# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_HandleSuspendKey_1.target
[Unit]
Description=Foo
[Login]
HandleSuspendKey=1

