# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_NAutoVTs_-1.target
[Unit]
Description=Foo
[Login]
NAutoVTs=-1

