# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_InhibitDelayMaxSec_True.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=True

