# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_IdleAction_True.target
[Unit]
Description=Foo
[Login]
IdleAction=True

