# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_HandleLidSwitch_True.target
[Unit]
Description=Foo
[Login]
HandleLidSwitch=True

