# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_KillUserProcesses_bar.target
[Unit]
Description=Foo
[Login]
KillUserProcesses=bar

