# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_HandleHibernateKey_1.target
[Unit]
Description=Foo
[Login]
HandleHibernateKey=1

