# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_IdleActionSec_3decades.target
[Unit]
Description=Foo
[Login]
IdleActionSec=3decades

