# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_HandleLidSwitchDocked_True.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchDocked=True

