# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_HandleLidSwitchExternalPower_True.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchExternalPower=True

