# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_InhibitorsMax_10000000000000T.target
[Unit]
Description=Foo
[Login]
InhibitorsMax=10000000000000T

