# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_IdleActionSec_0.5d.target
[Unit]
Description=Foo
[Login]
IdleActionSec=0.5d

