# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_HandleLidSwitchDocked_1.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchDocked=1

