# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_InhibitDelayMaxSec_0.5d.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=0.5d

