# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_RemoveIPC_foo.target
[Unit]
Description=Foo
[Login]
RemoveIPC=foo

