# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_InhibitDelayMaxSec_2potatoes.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=2potatoes

