# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_UserStopDelaySec_2potatoes.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=2potatoes

