# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_SuspendKeyIgnoreInhibited_bar.target
[Unit]
Description=Foo
[Login]
SuspendKeyIgnoreInhibited=bar

