# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_SessionsMax_-1T.target
[Unit]
Description=Foo
[Login]
SessionsMax=-1T

