# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_HandleHibernateKey_True.target
[Unit]
Description=Foo
[Login]
HandleHibernateKey=True

