# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_IdleActionSec_2potatoes.target
[Unit]
Description=Foo
[Login]
IdleActionSec=2potatoes

