# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_HandleLidSwitchExternalPower_1.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchExternalPower=1

