# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_InhibitDelayMaxSec_1p.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=1p

