# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_KillUserProcesses_baz.target
[Unit]
Description=Foo
[Login]
KillUserProcesses=baz

