# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_NAutoVTs_10000000000000.target
[Unit]
Description=Foo
[Login]
NAutoVTs=10000000000000

