# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_InhibitDelayMaxSec_3decades.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=3decades

