# systemdlint/tests/auto/InvalidValue/2.31/bad/Login_HandleLidSwitch_1.target
[Unit]
Description=Foo
[Login]
HandleLidSwitch=1

