# systemdlint/tests/auto/InvalidValue/2.31/bad/MACVTAP_Mode_True.target
[Unit]
Description=Foo
[MACVTAP]
Mode=True

