# systemdlint/tests/auto/InvalidValue/2.31/bad/MACVLAN_Mode_1.target
[Unit]
Description=Foo
[MACVLAN]
Mode=1

