# systemdlint/tests/auto/InvalidValue/2.31/bad/MACVLAN_Mode_True.target
[Unit]
Description=Foo
[MACVLAN]
Mode=True

