# systemdlint/tests/auto/InvalidValue/2.31/bad/MACVTAP_Mode_1.target
[Unit]
Description=Foo
[MACVTAP]
Mode=1

