# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DumpCore_baz.target
[Unit]
Description=Foo
[Manager]
DumpCore=baz

