# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultStartLimitBurst_True.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=True

