# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultStartLimitBurst_3decades.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=3decades

