# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultStartLimitBurst_0.5d.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=0.5d

