# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_ShutdownWatchdogSec_True.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=True

