# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultTimeoutStopSec_10000000000000.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=10000000000000

