# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_TimerSlackNSec_0.5d.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=0.5d

