# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultRestartSec_-1.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=-1

