# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultIPAccounting_foo.target
[Unit]
Description=Foo
[Manager]
DefaultIPAccounting=foo

