# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultEnvironment_True.target
[Unit]
Description=Foo
[Manager]
DefaultEnvironment=True

