# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultMemoryAccounting_baz.target
[Unit]
Description=Foo
[Manager]
DefaultMemoryAccounting=baz

