# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultTimeoutStopSec_2potatoes.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=2potatoes

