# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_CrashShell_foo.target
[Unit]
Description=Foo
[Manager]
CrashShell=foo

