# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_CPUAffinity_-1.target
[Unit]
Description=Foo
[Manager]
CPUAffinity=-1

