# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultIOAccounting_foo.target
[Unit]
Description=Foo
[Manager]
DefaultIOAccounting=foo

