# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultRestartSec_10000000000000.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=10000000000000

