# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultTasksMax_-1%25.target
[Unit]
Description=Foo
[Manager]
DefaultTasksMax=-1%

