# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_CrashReboot_bar.target
[Unit]
Description=Foo
[Manager]
CrashReboot=bar

