# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultTimeoutStopSec_True.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=True

