# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_LogColor_baz.target
[Unit]
Description=Foo
[Manager]
LogColor=baz

