# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultStartLimitIntervalSec_10000000000000.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=10000000000000

