# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DumpCore_foo.target
[Unit]
Description=Foo
[Manager]
DumpCore=foo

