# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_CrashReboot_baz.target
[Unit]
Description=Foo
[Manager]
CrashReboot=baz

