# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_RuntimeWatchdogSec_1p.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=1p

