# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultIPAccounting_baz.target
[Unit]
Description=Foo
[Manager]
DefaultIPAccounting=baz

