# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultStartLimitBurst_10000000000000.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=10000000000000

