# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultCPUAccounting_foo.target
[Unit]
Description=Foo
[Manager]
DefaultCPUAccounting=foo

