# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultMemoryAccounting_bar.target
[Unit]
Description=Foo
[Manager]
DefaultMemoryAccounting=bar

