# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_CrashChangeVT_foo.target
[Unit]
Description=Foo
[Manager]
CrashChangeVT=foo

