# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultStartLimitIntervalSec_1p.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=1p

