# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_TimerSlackNSec_-1.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=-1

