# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_WatchdogDevice_--devvvv.target
[Unit]
Description=Foo
[Manager]
WatchdogDevice=/devvvv

