# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultStartLimitIntervalSec_3decades.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=3decades

