# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultTimerAccuracySec_10000000000000.target
[Unit]
Description=Foo
[Manager]
DefaultTimerAccuracySec=10000000000000

