# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_TimerSlackNSec_10000000000000.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=10000000000000

