# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultTimeoutStopSec_-1.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=-1

