# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultTimeoutStartSec_10000000000000.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStartSec=10000000000000

