# systemdlint/tests/auto/InvalidValue/2.31/bad/Manager_DefaultTimeoutStartSec_True.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStartSec=True

