# systemdlint/tests/auto/InvalidValue/2.31/bad/Match_MACAddress_True.target
[Unit]
Description=Foo
[Match]
MACAddress=True

