# systemdlint/tests/auto/InvalidValue/2.31/bad/Match_MACAddress_AABB%3ACCDD%3AEEFF.target
[Unit]
Description=Foo
[Match]
MACAddress=AABB:CCDD:EEFF

