# systemdlint/tests/auto/InvalidValue/2.31/bad/Mount_DirectoryMode_888.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
DirectoryMode=888

