# systemdlint/tests/auto/InvalidValue/2.31/bad/Mount_RootImage_ftp%3A----bar.gov.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
RootImage=ftp://bar.gov

