# systemdlint/tests/auto/InvalidValue/2.31/bad/Mount_LimitCORE_-1.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
LimitCORE=-1

