# systemdlint/tests/auto/InvalidValue/2.31/bad/Mount_LimitRTTIME_10000000000000.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
LimitRTTIME=10000000000000

