# systemdlint/tests/auto/InvalidValue/2.31/bad/Mount_RootDirectory_ftp%3A----bar.gov.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
RootDirectory=ftp://bar.gov

