# systemdlint/tests/auto/InvalidValue/2.31/bad/Mount_WorkingDirectory_ftp%3A----bar.gov.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
WorkingDirectory=ftp://bar.gov

