# systemdlint/tests/auto/InvalidValue/2.31/bad/Mount_Group_4.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
Group=4

