# systemdlint/tests/auto/InvalidValue/2.31/bad/Mount_WorkingDirectory_http%3A----foo.org.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
WorkingDirectory=http://foo.org

