# systemdlint/tests/auto/InvalidValue/2.31/bad/Mount_TimeoutSec_0.5d.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
TimeoutSec=0.5d

