# systemdlint/tests/auto/InvalidValue/2.31/bad/Mount_BlockIOReadBandwidthMax_--devvvv%20-1M.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
BlockIOReadBandwidthMax=/devvvv -1M

