# systemdlint/tests/auto/InvalidValue/2.31/bad/Mount_MemoryMax_-1B.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
MemoryMax=-1B

