# systemdlint/tests/auto/InvalidValue/2.31/bad/Mount_BlockIOReadBandwidthMax_--tmp%2010000000000000K.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
BlockIOReadBandwidthMax=/tmp 10000000000000K

