# systemdlint/tests/auto/InvalidValue/2.31/bad/Mount_PrivateNetwork_baz.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
PrivateNetwork=baz

