# systemdlint/tests/auto/InvalidValue/2.31/bad/Neighbor_Address_0.0.0..target
[Unit]
Description=Foo
[Neighbor]
Address=0.0.0.

