# systemdlint/tests/auto/InvalidValue/2.31/bad/Neighbor_Address_255.255.255.256.target
[Unit]
Description=Foo
[Neighbor]
Address=255.255.255.256

