# systemdlint/tests/auto/InvalidValue/2.31/bad/NetDev_MTUBytes_9217.target
[Unit]
Description=Foo
[NetDev]
Name=abc
MTUBytes=9217

