# systemdlint/tests/auto/InvalidValue/2.31/bad/NetDev_MACAddress_00-11-22-GG-44-55.target
[Unit]
Description=Foo
[NetDev]
Name=abc
MACAddress=00-11-22-GG-44-55

