# systemdlint/tests/auto/InvalidValue/2.31/bad/NetDev_MACAddress_AABB%3ACCDD%3AEEFF.target
[Unit]
Description=Foo
[NetDev]
Name=abc
MACAddress=AABB:CCDD:EEFF

