# systemdlint/tests/auto/InvalidValue/2.31/bad/NetDev_Kind_1.target
[Unit]
Description=Foo
[NetDev]
Name=abc
Kind=1

