# systemdlint/tests/auto/InvalidValue/2.31/bad/NetDev_MACAddress_True.target
[Unit]
Description=Foo
[NetDev]
Name=abc
MACAddress=True

