# systemdlint/tests/auto/InvalidValue/2.31/bad/NetDev_MTUBytes_5121.target
[Unit]
Description=Foo
[NetDev]
Name=abc
MTUBytes=5121

