# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_IgnoreCarrierLoss_bar.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=bar

