# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_LLDP_True.target
[Unit]
Description=Foo
[Network]
LLDP=True

