# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_PrimarySlave_bar.target
[Unit]
Description=Foo
[Network]
PrimarySlave=bar

