# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_MulticastDNS_True.target
[Unit]
Description=Foo
[Network]
MulticastDNS=True

