# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_NTP_http%3Atest.target
[Unit]
Description=Foo
[Network]
NTP=http:test

