# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_PrimarySlave_foo.target
[Unit]
Description=Foo
[Network]
PrimarySlave=foo

