# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_PrimarySlave_baz.target
[Unit]
Description=Foo
[Network]
PrimarySlave=baz

