# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_Gateway_256.255.255.255.target
[Unit]
Description=Foo
[Network]
Gateway=256.255.255.255

