# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_IPv6Token_256.255.255.255.target
[Unit]
Description=Foo
[Network]
IPv6Token=256.255.255.255

