# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_LinkLocalAddressing_True.target
[Unit]
Description=Foo
[Network]
LinkLocalAddressing=True

