# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_ActiveSlave_foo.target
[Unit]
Description=Foo
[Network]
ActiveSlave=foo

