# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_IPv6MTUBytes_10000000000000.target
[Unit]
Description=Foo
[Network]
IPv6MTUBytes=10000000000000

