# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_DNS_255.256.255.255.target
[Unit]
Description=Foo
[Network]
DNS=255.256.255.255

