# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_IPMasquerade_bar.target
[Unit]
Description=Foo
[Network]
IPMasquerade=bar

