# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_IPv6HopLimit_10000000000000.target
[Unit]
Description=Foo
[Network]
IPv6HopLimit=10000000000000

