# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_DNSDefaultRoute_foo.target
[Unit]
Description=Foo
[Network]
DNSDefaultRoute=foo

