# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_IgnoreCarrierLoss_baz.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=baz

