# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_NTP_abc.target
[Unit]
Description=Foo
[Network]
NTP=abc

