# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_DNSOverTLS_1.target
[Unit]
Description=Foo
[Network]
DNSOverTLS=1

