# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_DNS_0--abc.target
[Unit]
Description=Foo
[Network]
DNS=0/abc

