# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_Address_255.256.255.255.target
[Unit]
Description=Foo
[Network]
Address=255.256.255.255

