# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_Gateway_abc.target
[Unit]
Description=Foo
[Network]
Gateway=abc

