# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_DHCPServer_bar.target
[Unit]
Description=Foo
[Network]
DHCPServer=bar

