# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_DNSOverTLS_True.target
[Unit]
Description=Foo
[Network]
DNSOverTLS=True

