# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_LLMNR_True.target
[Unit]
Description=Foo
[Network]
LLMNR=True

