# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_IPv6ProxyNDP_foo.target
[Unit]
Description=Foo
[Network]
IPv6ProxyNDP=foo

