# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_Gateway_1.1.1.1---1.target
[Unit]
Description=Foo
[Network]
Gateway=1.1.1.1/-1

