# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_DNSSEC_1.target
[Unit]
Description=Foo
[Network]
DNSSEC=1

