# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_NTP_0--abc.target
[Unit]
Description=Foo
[Network]
NTP=0/abc

