# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_ConfigureWithoutCarrier_foo.target
[Unit]
Description=Foo
[Network]
ConfigureWithoutCarrier=foo

