# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_LinkLocalAddressing_1.target
[Unit]
Description=Foo
[Network]
LinkLocalAddressing=1

