# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_NTP_True.target
[Unit]
Description=Foo
[Network]
NTP=True

