# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_IPv6ProxyNDPAddress_256.255.255.255.target
[Unit]
Description=Foo
[Network]
IPv6ProxyNDPAddress=256.255.255.255

