# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_IPv6Token_-1.1.1.1.target
[Unit]
Description=Foo
[Network]
IPv6Token=-1.1.1.1

