# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_Address_0--abc.target
[Unit]
Description=Foo
[Network]
Address=0/abc

