# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_NTP_1.target
[Unit]
Description=Foo
[Network]
NTP=1

