# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_EmitLLDP_1.target
[Unit]
Description=Foo
[Network]
EmitLLDP=1

