# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_IPv4LLRoute_baz.target
[Unit]
Description=Foo
[Network]
IPv4LLRoute=baz

