# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_DNSSEC_True.target
[Unit]
Description=Foo
[Network]
DNSSEC=True

