# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_DHCP_1.target
[Unit]
Description=Foo
[Network]
DHCP=1

