# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_Gateway_0--abc.target
[Unit]
Description=Foo
[Network]
Gateway=0/abc

