# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_NTP_256.255.255.255.target
[Unit]
Description=Foo
[Network]
NTP=256.255.255.255

