# systemdlint/tests/auto/InvalidValue/2.31/bad/Network_DHCP_True.target
[Unit]
Description=Foo
[Network]
DHCP=True

