# systemdlint/tests/auto/InvalidValue/2.31/bad/Path_DirectoryMode_-1.target
[Unit]
Description=Foo
[Path]
DirectoryMode=-1

