# systemdlint/tests/auto/InvalidValue/2.31/bad/Path_PathChanged_tmp.target
[Unit]
Description=Foo
[Path]
PathChanged=tmp

