# systemdlint/tests/auto/InvalidValue/2.31/bad/Path_DirectoryNotEmpty_1.target
[Unit]
Description=Foo
[Path]
DirectoryNotEmpty=1

