# systemdlint/tests/auto/InvalidValue/2.31/bad/Path_PathChanged_..--..--.target
[Unit]
Description=Foo
[Path]
PathChanged=../../

