# systemdlint/tests/auto/InvalidValue/2.31/bad/Path_DirectoryMode_999.target
[Unit]
Description=Foo
[Path]
DirectoryMode=999

