# systemdlint/tests/auto/InvalidValue/2.31/bad/Path_PathModified_True.target
[Unit]
Description=Foo
[Path]
PathModified=True

