# systemdlint/tests/auto/InvalidValue/2.31/bad/Path_MakeDirectory_bar.target
[Unit]
Description=Foo
[Path]
MakeDirectory=bar

