# systemdlint/tests/auto/InvalidValue/2.31/bad/Path_MakeDirectory_baz.target
[Unit]
Description=Foo
[Path]
MakeDirectory=baz

