# systemdlint/tests/auto/InvalidValue/2.31/bad/Path_PathChanged_1.target
[Unit]
Description=Foo
[Path]
PathChanged=1

