# systemdlint/tests/auto/InvalidValue/2.31/bad/Peer_MACAddress_AABB%3ACCDD%3AEEFF.target
[Unit]
Description=Foo
[Peer]
Name=abc
MACAddress=AABB:CCDD:EEFF

