# systemdlint/tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_TypeOfService_-1.target
[Unit]
Description=Foo
[RoutingPolicyRule]
TypeOfService=-1

