# systemdlint/tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_InvertRule_bar.target
[Unit]
Description=Foo
[RoutingPolicyRule]
InvertRule=bar

