# systemdlint/tests/auto/InvalidValue/2.31/bad/Route_Gateway_0--abc.target
[Unit]
Description=Foo
[Route]
Gateway=0/abc

