# systemdlint/tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_SourcePort_-1%20-1.target
[Unit]
Description=Foo
[RoutingPolicyRule]
SourcePort=-1 -1

