# systemdlint/tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_DestinationPort_10000000000000%20-1.target
[Unit]
Description=Foo
[RoutingPolicyRule]
DestinationPort=10000000000000 -1

