# systemdlint/tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_IPProtocol_1.target
[Unit]
Description=Foo
[RoutingPolicyRule]
IPProtocol=1

