# systemdlint/tests/auto/InvalidValue/2.31/bad/Route_Gateway_1.1.1.1---1.target
[Unit]
Description=Foo
[Route]
Gateway=1.1.1.1/-1

