# systemdlint/tests/auto/InvalidValue/2.31/bad/Route_Metric_-1.target
[Unit]
Description=Foo
[Route]
Metric=-1

