# systemdlint/tests/auto/InvalidValue/2.31/bad/Route_GatewayOnLink_foo.target
[Unit]
Description=Foo
[Route]
GatewayOnLink=foo

