# systemdlint/tests/auto/InvalidValue/2.31/bad/Route_Gateway_abc.target
[Unit]
Description=Foo
[Route]
Gateway=abc

