# systemdlint/tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_IPProtocol_True.target
[Unit]
Description=Foo
[RoutingPolicyRule]
IPProtocol=True

