# systemdlint/tests/auto/InvalidValue/2.31/bad/Route_Table_0.target
[Unit]
Description=Foo
[Route]
Table=0

