# systemdlint/tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_From_255.256.255.255.target
[Unit]
Description=Foo
[RoutingPolicyRule]
From=255.256.255.255

