# systemdlint/tests/auto/InvalidValue/2.31/bad/Route_Scope_True.target
[Unit]
Description=Foo
[Route]
Scope=True

