# systemdlint/tests/auto/InvalidValue/2.31/bad/Route_QuickAck_foo.target
[Unit]
Description=Foo
[Route]
QuickAck=foo

