# systemdlint/tests/auto/InvalidValue/2.31/bad/Route_Type_True.target
[Unit]
Description=Foo
[Route]
Type=True

