# systemdlint/tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_InvertRule_foo.target
[Unit]
Description=Foo
[RoutingPolicyRule]
InvertRule=foo

