# systemdlint/tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_TypeOfService_256.target
[Unit]
Description=Foo
[RoutingPolicyRule]
TypeOfService=256

