# systemdlint/tests/auto/InvalidValue/2.31/bad/Route_IPv6Preference_1.target
[Unit]
Description=Foo
[Route]
IPv6Preference=1

