# systemdlint/tests/auto/InvalidValue/2.31/bad/Route_Scope_1.target
[Unit]
Description=Foo
[Route]
Scope=1

