# systemdlint/tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_From_0--abc.target
[Unit]
Description=Foo
[RoutingPolicyRule]
From=0/abc

