# systemdlint/tests/auto/InvalidValue/2.31/bad/Route_IPv6Preference_True.target
[Unit]
Description=Foo
[Route]
IPv6Preference=True

