# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStartPre_--etc--passwd.target
[Unit]
Description=Foo
[Service]
ExecStartPre=/etc/passwd

