# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_BlockIOReadBandwidthMax_dev--foo%2010000000000000K.target
[Unit]
Description=Foo
[Service]
BlockIOReadBandwidthMax=dev/foo 10000000000000K

