# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_Group_1.target
[Unit]
Description=Foo
[Service]
Group=1

