# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_BlockIOReadBandwidthMax_--devvvv%20-1K.target
[Unit]
Description=Foo
[Service]
BlockIOReadBandwidthMax=/devvvv -1K

