# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_SyslogLevelPrefix_baz.target
[Unit]
Description=Foo
[Service]
SyslogLevelPrefix=baz

