# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_GuessMainPID_foo.target
[Unit]
Description=Foo
[Service]
GuessMainPID=foo

