# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_WatchdogSec_True.target
[Unit]
Description=Foo
[Service]
WatchdogSec=True

