# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_DevicePolicy_True.target
[Unit]
Description=Foo
[Service]
DevicePolicy=True

