# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_TasksAccounting_bar.target
[Unit]
Description=Foo
[Service]
TasksAccounting=bar

