# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_SuccessExitStatus_True.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=True

