# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_UtmpMode_1.target
[Unit]
Description=Foo
[Service]
UtmpMode=1

