# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_NotifyAccess_1.target
[Unit]
Description=Foo
[Service]
NotifyAccess=1

