# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_Personality_1.target
[Unit]
Description=Foo
[Service]
Personality=1

