# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_GuessMainPID_bar.target
[Unit]
Description=Foo
[Service]
GuessMainPID=bar

