# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_BlockIODeviceWeight_--devvvv%2010001.target
[Unit]
Description=Foo
[Service]
BlockIODeviceWeight=/devvvv 10001

