# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_SyslogFacility_1.target
[Unit]
Description=Foo
[Service]
SyslogFacility=1

