# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_Group_True.target
[Unit]
Description=Foo
[Service]
Group=True

