# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_BlockIOReadBandwidthMax_--tmp%2010000000000000K.target
[Unit]
Description=Foo
[Service]
BlockIOReadBandwidthMax=/tmp 10000000000000K

