# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_EnvironmentFile_1.target
[Unit]
Description=Foo
[Service]
EnvironmentFile=1

