# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_MountFlags_True.target
[Unit]
Description=Foo
[Service]
MountFlags=True

