# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStart_1.target
[Unit]
Description=Foo
[Service]
ExecStart=1

