# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_LimitSTACK_10000000000000.target
[Unit]
Description=Foo
[Service]
LimitSTACK=10000000000000

