# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_TimerSlackNSec_3decades.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=3decades

