# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_RootDirectoryStartOnly_baz.target
[Unit]
Description=Foo
[Service]
RootDirectoryStartOnly=baz

