# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_RestrictSUIDSGID_bar.target
[Unit]
Description=Foo
[Service]
RestrictSUIDSGID=bar

