# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_Type_1.target
[Unit]
Description=Foo
[Service]
Type=1

