# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_LimitAS_10000000000000.target
[Unit]
Description=Foo
[Service]
LimitAS=10000000000000

