# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_RestartPreventExitStatus_True.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=True

