# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_IgnoreSIGPIPE_bar.target
[Unit]
Description=Foo
[Service]
IgnoreSIGPIPE=bar

