# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_KillMode_1.target
[Unit]
Description=Foo
[Service]
KillMode=1

