# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_RootDirectory_http%3A----foo.org.target
[Unit]
Description=Foo
[Service]
RootDirectory=http://foo.org

