# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_LimitCPU_10000000000000.target
[Unit]
Description=Foo
[Service]
LimitCPU=10000000000000

