# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecReload_1.target
[Unit]
Description=Foo
[Service]
ExecReload=1

