# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_Personality_True.target
[Unit]
Description=Foo
[Service]
Personality=True

