# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_WatchdogSignal_True.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=True

