# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_RemoveIPC_bar.target
[Unit]
Description=Foo
[Service]
RemoveIPC=bar

