# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_WorkingDirectory_http%3A----foo.org.target
[Unit]
Description=Foo
[Service]
WorkingDirectory=http://foo.org

