# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_BlockIODeviceWeight_--tmp--dev%200.target
[Unit]
Description=Foo
[Service]
BlockIODeviceWeight=/tmp/dev 0

