# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_TimeoutSec_True.target
[Unit]
Description=Foo
[Service]
TimeoutSec=True

