# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_RestartForceExitStatus_True.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=True

