# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_User_4.target
[Unit]
Description=Foo
[Service]
User=4

