# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStop_True.target
[Unit]
Description=Foo
[Service]
ExecStop=True

