# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_SyslogLevel_1.target
[Unit]
Description=Foo
[Service]
SyslogLevel=1

