# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_RuntimeMaxSec_3decades.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=3decades

