# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_StartLimitInterval_True.target
[Unit]
Description=Foo
[Service]
StartLimitInterval=True

