# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_MountAPIVFS_bar.target
[Unit]
Description=Foo
[Service]
MountAPIVFS=bar

