# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStopPost_--etc--passwd.target
[Unit]
Description=Foo
[Service]
ExecStopPost=/etc/passwd

