# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_RestrictSUIDSGID_foo.target
[Unit]
Description=Foo
[Service]
RestrictSUIDSGID=foo

