# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_CPUSchedulingResetOnFork_foo.target
[Unit]
Description=Foo
[Service]
CPUSchedulingResetOnFork=foo

