# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_SyslogFacility_True.target
[Unit]
Description=Foo
[Service]
SyslogFacility=True

