# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStartPre_True.target
[Unit]
Description=Foo
[Service]
ExecStartPre=True

