# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_MountAPIVFS_foo.target
[Unit]
Description=Foo
[Service]
MountAPIVFS=foo

