# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_CPUSchedulingPolicy_1.target
[Unit]
Description=Foo
[Service]
CPUSchedulingPolicy=1

