# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_WorkingDirectory_ftp%3A----bar.gov.target
[Unit]
Description=Foo
[Service]
WorkingDirectory=ftp://bar.gov

