# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_PrivateTmp_bar.target
[Unit]
Description=Foo
[Service]
PrivateTmp=bar

