# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_DevicePolicy_1.target
[Unit]
Description=Foo
[Service]
DevicePolicy=1

