# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_BlockIOWriteBandwidthMax_dev--foo%20-1T.target
[Unit]
Description=Foo
[Service]
BlockIOWriteBandwidthMax=dev/foo -1T

