# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_SupplementaryGroups_True.target
[Unit]
Description=Foo
[Service]
SupplementaryGroups=True

