# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_StartLimitIntervalSec_-1.target
[Unit]
Description=Foo
[Service]
StartLimitIntervalSec=-1

