# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_TimeoutStopSec_0.5d.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=0.5d

