# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_LimitNPROC_10000000000000.target
[Unit]
Description=Foo
[Service]
LimitNPROC=10000000000000

