# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_SyslogLevel_True.target
[Unit]
Description=Foo
[Service]
SyslogLevel=True

