# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_MountAPIVFS_baz.target
[Unit]
Description=Foo
[Service]
MountAPIVFS=baz

