# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_IOSchedulingClass_True.target
[Unit]
Description=Foo
[Service]
IOSchedulingClass=True

