# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStartPost_1.target
[Unit]
Description=Foo
[Service]
ExecStartPost=1

