# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStart_--etc--passwd.target
[Unit]
Description=Foo
[Service]
ExecStart=/etc/passwd

