# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_EnvironmentFile_..--..--.target
[Unit]
Description=Foo
[Service]
EnvironmentFile=../../

