# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_RootDirectory_ftp%3A----bar.gov.target
[Unit]
Description=Foo
[Service]
RootDirectory=ftp://bar.gov

