# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStartPost_True.target
[Unit]
Description=Foo
[Service]
ExecStartPost=True

