# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_SupplementaryGroups_1.target
[Unit]
Description=Foo
[Service]
SupplementaryGroups=1

