# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_NotifyAccess_True.target
[Unit]
Description=Foo
[Service]
NotifyAccess=True

