# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStart_True.target
[Unit]
Description=Foo
[Service]
ExecStart=True

