# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStartPost_--etc--passwd.target
[Unit]
Description=Foo
[Service]
ExecStartPost=/etc/passwd

