# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_User_meisnotausercozofmyloooooooooooooooooooooooooooooooooognamegee.target
[Unit]
Description=Foo
[Service]
User=meisnotausercozofmyloooooooooooooooooooooooooooooooooognamegee

