# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_TasksAccounting_foo.target
[Unit]
Description=Foo
[Service]
TasksAccounting=foo

