# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStart_a%3Eb.target
[Unit]
Description=Foo
[Service]
ExecStart=a>b

