# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_RestrictRealtime_foo.target
[Unit]
Description=Foo
[Service]
RestrictRealtime=foo

