# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_TimerAccuracySec_0.5d.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=0.5d

