# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_RemoveIPC_foo.target
[Unit]
Description=Foo
[Service]
RemoveIPC=foo

