# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_NonBlocking_baz.target
[Unit]
Description=Foo
[Service]
NonBlocking=baz

