# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_TTYVHangup_bar.target
[Unit]
Description=Foo
[Service]
TTYVHangup=bar

