# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_RemoveIPC_baz.target
[Unit]
Description=Foo
[Service]
RemoveIPC=baz

