# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecReload_a%3Eb.target
[Unit]
Description=Foo
[Service]
ExecReload=a>b

