# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_DeviceAllow_--devvvv.target
[Unit]
Description=Foo
[Service]
DeviceAllow=/devvvv

