# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_PIDFile_ftp%3A----bar.gov.target
[Unit]
Description=Foo
[Service]
PIDFile=ftp://bar.gov

