# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_LimitSTACK_-1.target
[Unit]
Description=Foo
[Service]
LimitSTACK=-1

