# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_LimitCORE_10000000000000.target
[Unit]
Description=Foo
[Service]
LimitCORE=10000000000000

