# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_MemoryLow_-1B.target
[Unit]
Description=Foo
[Service]
MemoryLow=-1B

