# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_TimerSlackNSec_2potatoes.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=2potatoes

