# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStop_1.target
[Unit]
Description=Foo
[Service]
ExecStop=1

