# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_Restart_True.target
[Unit]
Description=Foo
[Service]
Restart=True

