# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_TTYPath_--tmp--dev.target
[Unit]
Description=Foo
[Service]
TTYPath=/tmp/dev

