# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_WatchdogSec_3decades.target
[Unit]
Description=Foo
[Service]
WatchdogSec=3decades

