# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_LimitRTPRIO_10000000000000.target
[Unit]
Description=Foo
[Service]
LimitRTPRIO=10000000000000

