# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_TimeoutStopSec_3decades.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=3decades

