# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_TimerSlackNSec_10000000000000.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=10000000000000

