# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_PIDFile_http%3A----foo.org.target
[Unit]
Description=Foo
[Service]
PIDFile=http://foo.org

