# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_LimitDATA_10000000000000.target
[Unit]
Description=Foo
[Service]
LimitDATA=10000000000000

