# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_LimitNICE_-1.target
[Unit]
Description=Foo
[Service]
LimitNICE=-1

