# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_Delegate_1.target
[Unit]
Description=Foo
[Service]
Delegate=1

