# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStopPost_a%3Eb.target
[Unit]
Description=Foo
[Service]
ExecStopPost=a>b

