# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_SupplementaryGroups_meisnotagroupcozofmyloooooooooooooooooooooooooooooooooognamegee.target
[Unit]
Description=Foo
[Service]
SupplementaryGroups=meisnotagroupcozofmyloooooooooooooooooooooooooooooooooognamegee

