# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_TimeoutStartSec_True.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=True

