# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_Type_True.target
[Unit]
Description=Foo
[Service]
Type=True

