# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_FileDescriptorStoreMax_10000000000000.target
[Unit]
Description=Foo
[Service]
FileDescriptorStoreMax=10000000000000

