# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecReload_True.target
[Unit]
Description=Foo
[Service]
ExecReload=True

