# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_RootImage_http%3A----foo.org.target
[Unit]
Description=Foo
[Service]
RootImage=http://foo.org

