# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_CPUSchedulingPolicy_True.target
[Unit]
Description=Foo
[Service]
CPUSchedulingPolicy=True

