# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_RootImage_ftp%3A----bar.gov.target
[Unit]
Description=Foo
[Service]
RootImage=ftp://bar.gov

