# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_StartLimitIntervalSec_10000000000000.target
[Unit]
Description=Foo
[Service]
StartLimitIntervalSec=10000000000000

