# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStartPre_a%3Eb.target
[Unit]
Description=Foo
[Service]
ExecStartPre=a>b

