# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_TTYPath_--devvvv.target
[Unit]
Description=Foo
[Service]
TTYPath=/devvvv

