# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_SecureBits_1.target
[Unit]
Description=Foo
[Service]
SecureBits=1

