# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStopPost_True.target
[Unit]
Description=Foo
[Service]
ExecStopPost=True

