# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_RestartSec_3decades.target
[Unit]
Description=Foo
[Service]
RestartSec=3decades

