# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_SyslogLevelPrefix_bar.target
[Unit]
Description=Foo
[Service]
SyslogLevelPrefix=bar

