# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStop_--etc--passwd.target
[Unit]
Description=Foo
[Service]
ExecStop=/etc/passwd

