# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_CPUSchedulingResetOnFork_baz.target
[Unit]
Description=Foo
[Service]
CPUSchedulingResetOnFork=baz

