# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_KillSignal_True.target
[Unit]
Description=Foo
[Service]
KillSignal=True

