# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_FinalKillSignal_True.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=True

