# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_TimerAccuracySec_True.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=True

