# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_Restart_1.target
[Unit]
Description=Foo
[Service]
Restart=1

