# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_BlockIOWriteBandwidthMax_--tmp--dev%20-1K.target
[Unit]
Description=Foo
[Service]
BlockIOWriteBandwidthMax=/tmp/dev -1K

