# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_UMask_10000000000000.target
[Unit]
Description=Foo
[Service]
UMask=10000000000000

