# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_BlockIOWeight_0.target
[Unit]
Description=Foo
[Service]
BlockIOWeight=0

