# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_TTYVHangup_foo.target
[Unit]
Description=Foo
[Service]
TTYVHangup=foo

