# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_MountFlags_1.target
[Unit]
Description=Foo
[Service]
MountFlags=1

