# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_MemoryDenyWriteExecute_baz.target
[Unit]
Description=Foo
[Service]
MemoryDenyWriteExecute=baz

