# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_MemoryDenyWriteExecute_bar.target
[Unit]
Description=Foo
[Service]
MemoryDenyWriteExecute=bar

