# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_TimerAccuracySec_3decades.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=3decades

