# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStop_a%3Eb.target
[Unit]
Description=Foo
[Service]
ExecStop=a>b

