# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_RestartSec_True.target
[Unit]
Description=Foo
[Service]
RestartSec=True

