# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_LimitAS_-1.target
[Unit]
Description=Foo
[Service]
LimitAS=-1

