# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStopPost_1.target
[Unit]
Description=Foo
[Service]
ExecStopPost=1

