# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_LimitSIGPENDING_10000000000000.target
[Unit]
Description=Foo
[Service]
LimitSIGPENDING=10000000000000

