# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_LimitMSGQUEUE_-1.target
[Unit]
Description=Foo
[Service]
LimitMSGQUEUE=-1

