# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_MemoryMax_10000000000000T.target
[Unit]
Description=Foo
[Service]
MemoryMax=10000000000000T

