# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_CPUSchedulingPriority_0.target
[Unit]
Description=Foo
[Service]
CPUSchedulingPriority=0

