# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_BlockIODeviceWeight_--tmp--dev%2010001.target
[Unit]
Description=Foo
[Service]
BlockIODeviceWeight=/tmp/dev 10001

