# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_SendSIGKILL_baz.target
[Unit]
Description=Foo
[Service]
SendSIGKILL=baz

