# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_IOSchedulingClass_1.target
[Unit]
Description=Foo
[Service]
IOSchedulingClass=1

