# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_CPUSchedulingPriority_100.target
[Unit]
Description=Foo
[Service]
CPUSchedulingPriority=100

