# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_ExecStartPost_a%3Eb.target
[Unit]
Description=Foo
[Service]
ExecStartPost=a>b

