# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_SecureBits_True.target
[Unit]
Description=Foo
[Service]
SecureBits=True

