# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_LimitRTPRIO_-1.target
[Unit]
Description=Foo
[Service]
LimitRTPRIO=-1

