# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_RestrictRealtime_baz.target
[Unit]
Description=Foo
[Service]
RestrictRealtime=baz

