# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_Delegate_True.target
[Unit]
Description=Foo
[Service]
Delegate=True

