# systemdlint/tests/auto/InvalidValue/2.31/bad/Service_TimerSlackNSec_0.5d.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=0.5d

