# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_CPUSchedulingPolicy_True.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingPolicy=True

