# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_RootImage_http%3A----foo.org.target
[Unit]
Description=Foo
[Socket]
RootImage=http://foo.org

