# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_BindIPv6Only_True.target
[Unit]
Description=Foo
[Socket]
BindIPv6Only=True

