# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SocketProtocol_1.target
[Unit]
Description=Foo
[Socket]
SocketProtocol=1

