# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_TTYReset_foo.target
[Unit]
Description=Foo
[Socket]
TTYReset=foo

