# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_KillSignal_True.target
[Unit]
Description=Foo
[Socket]
KillSignal=True

