# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ExecStartPost_1.target
[Unit]
Description=Foo
[Socket]
ExecStartPost=1

