# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_TTYVHangup_baz.target
[Unit]
Description=Foo
[Socket]
TTYVHangup=baz

