# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SELinuxContextFromNet_baz.target
[Unit]
Description=Foo
[Socket]
SELinuxContextFromNet=baz

