# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_UtmpMode_True.target
[Unit]
Description=Foo
[Socket]
UtmpMode=True

