# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ExecStartPost_a%3Eb.target
[Unit]
Description=Foo
[Socket]
ExecStartPost=a>b

