# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SecureBits_True.target
[Unit]
Description=Foo
[Socket]
SecureBits=True

