# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_LimitRSS_-1.target
[Unit]
Description=Foo
[Socket]
LimitRSS=-1

