# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_StartupCPUShares_262145.target
[Unit]
Description=Foo
[Socket]
StartupCPUShares=262145

