# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_Broadcast_foo.target
[Unit]
Description=Foo
[Socket]
Broadcast=foo

