# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_Delegate_True.target
[Unit]
Description=Foo
[Socket]
Delegate=True

