# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_Delegate_1.target
[Unit]
Description=Foo
[Socket]
Delegate=1

