# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ExecStartPost_--etc--passwd.target
[Unit]
Description=Foo
[Socket]
ExecStartPost=/etc/passwd

