# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_LimitNICE_-1.target
[Unit]
Description=Foo
[Socket]
LimitNICE=-1

