# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_MessageQueueMaxMessages_-1.target
[Unit]
Description=Foo
[Socket]
MessageQueueMaxMessages=-1

