# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_UMask_10000000000000.target
[Unit]
Description=Foo
[Socket]
UMask=10000000000000

