# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_TTYPath_--tmp.target
[Unit]
Description=Foo
[Socket]
TTYPath=/tmp

