# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ExecStopPost_True.target
[Unit]
Description=Foo
[Socket]
ExecStopPost=True

