# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_BlockIOReadBandwidthMax_--tmp--dev%20-1K.target
[Unit]
Description=Foo
[Socket]
BlockIOReadBandwidthMax=/tmp/dev -1K

