# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ReceiveBuffer_-1M.target
[Unit]
Description=Foo
[Socket]
ReceiveBuffer=-1M

