# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_BlockIOWeight_0.target
[Unit]
Description=Foo
[Socket]
BlockIOWeight=0

