# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_Writable_baz.target
[Unit]
Description=Foo
[Socket]
Writable=baz

