# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_MountAPIVFS_foo.target
[Unit]
Description=Foo
[Socket]
MountAPIVFS=foo

