# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_TTYVHangup_foo.target
[Unit]
Description=Foo
[Socket]
TTYVHangup=foo

