# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_StartupIOWeight_10001.target
[Unit]
Description=Foo
[Socket]
StartupIOWeight=10001

