# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_Backlog_10000000000000.target
[Unit]
Description=Foo
[Socket]
Backlog=10000000000000

