# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ExecStopPost_--etc--passwd.target
[Unit]
Description=Foo
[Socket]
ExecStopPost=/etc/passwd

