# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_MemoryMax_10000000000000K.target
[Unit]
Description=Foo
[Socket]
MemoryMax=10000000000000K

