# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_BlockIOWriteBandwidthMax_--tmp%20-1G.target
[Unit]
Description=Foo
[Socket]
BlockIOWriteBandwidthMax=/tmp -1G

