# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_IgnoreSIGPIPE_foo.target
[Unit]
Description=Foo
[Socket]
IgnoreSIGPIPE=foo

