# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_Environment_1.target
[Unit]
Description=Foo
[Socket]
Environment=1

