# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SecureBits_1.target
[Unit]
Description=Foo
[Socket]
SecureBits=1

