# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_DevicePolicy_1.target
[Unit]
Description=Foo
[Socket]
DevicePolicy=1

