# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_Accept_bar.target
[Unit]
Description=Foo
[Socket]
Accept=bar

