# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_PrivateTmp_baz.target
[Unit]
Description=Foo
[Socket]
PrivateTmp=baz

