# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_DeferAcceptSec_10000000000000.target
[Unit]
Description=Foo
[Socket]
DeferAcceptSec=10000000000000

