# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_Transparent_baz.target
[Unit]
Description=Foo
[Socket]
Transparent=baz

