# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_LimitNPROC_10000000000000.target
[Unit]
Description=Foo
[Socket]
LimitNPROC=10000000000000

