# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_LimitDATA_-1.target
[Unit]
Description=Foo
[Socket]
LimitDATA=-1

