# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ListenFIFO_1.target
[Unit]
Description=Foo
[Socket]
ListenFIFO=1

