# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_LimitRTPRIO_10000000000000.target
[Unit]
Description=Foo
[Socket]
LimitRTPRIO=10000000000000

