# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_BlockStartupIOWeight_10001.target
[Unit]
Description=Foo
[Socket]
BlockStartupIOWeight=10001

