# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_WorkingDirectory_http%3A----foo.org.target
[Unit]
Description=Foo
[Socket]
WorkingDirectory=http://foo.org

