# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_EnvironmentFile_1.target
[Unit]
Description=Foo
[Socket]
EnvironmentFile=1

