# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SendSIGKILL_foo.target
[Unit]
Description=Foo
[Socket]
SendSIGKILL=foo

