# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SendSIGHUP_bar.target
[Unit]
Description=Foo
[Socket]
SendSIGHUP=bar

