# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_MemoryDenyWriteExecute_baz.target
[Unit]
Description=Foo
[Socket]
MemoryDenyWriteExecute=baz

