# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SystemCallErrorNumber_4096.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=4096

