# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SendSIGHUP_foo.target
[Unit]
Description=Foo
[Socket]
SendSIGHUP=foo

