# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_LimitRTPRIO_-1.target
[Unit]
Description=Foo
[Socket]
LimitRTPRIO=-1

