# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_MemoryHigh_-1G.target
[Unit]
Description=Foo
[Socket]
MemoryHigh=-1G

