# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_WatchdogSignal_True.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=True

