# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_LimitRTTIME_10000000000000.target
[Unit]
Description=Foo
[Socket]
LimitRTTIME=10000000000000

