# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_DisableControllers_1.target
[Unit]
Description=Foo
[Socket]
DisableControllers=1

