# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_DeferAcceptSec_-1.target
[Unit]
Description=Foo
[Socket]
DeferAcceptSec=-1

