# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SyslogLevelPrefix_baz.target
[Unit]
Description=Foo
[Socket]
SyslogLevelPrefix=baz

