# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_KeepAliveIntervalSec_-1.target
[Unit]
Description=Foo
[Socket]
KeepAliveIntervalSec=-1

