# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SocketGroup_3.target
[Unit]
Description=Foo
[Socket]
SocketGroup=3

