# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_Group_True.target
[Unit]
Description=Foo
[Socket]
Group=True

