# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ExecStopPre_a%3Eb.target
[Unit]
Description=Foo
[Socket]
ExecStopPre=a>b

