# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_TTYVHangup_bar.target
[Unit]
Description=Foo
[Socket]
TTYVHangup=bar

