# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_RootDirectory_ftp%3A----bar.gov.target
[Unit]
Description=Foo
[Socket]
RootDirectory=ftp://bar.gov

