# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_Writable_bar.target
[Unit]
Description=Foo
[Socket]
Writable=bar

