# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SyslogFacility_1.target
[Unit]
Description=Foo
[Socket]
SyslogFacility=1

