# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_PassSecurity_foo.target
[Unit]
Description=Foo
[Socket]
PassSecurity=foo

