# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_PrivateDevices_foo.target
[Unit]
Description=Foo
[Socket]
PrivateDevices=foo

