# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_RootDirectory_http%3A----foo.org.target
[Unit]
Description=Foo
[Socket]
RootDirectory=http://foo.org

