# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SyslogLevel_True.target
[Unit]
Description=Foo
[Socket]
SyslogLevel=True

