# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_TasksAccounting_foo.target
[Unit]
Description=Foo
[Socket]
TasksAccounting=foo

