# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_MemoryAccounting_baz.target
[Unit]
Description=Foo
[Socket]
MemoryAccounting=baz

