# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_CPUSchedulingPolicy_1.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingPolicy=1

