# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_DevicePolicy_True.target
[Unit]
Description=Foo
[Socket]
DevicePolicy=True

