# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_BlockIOWriteBandwidthMax_--tmp--dev%2010000000000000T.target
[Unit]
Description=Foo
[Socket]
BlockIOWriteBandwidthMax=/tmp/dev 10000000000000T

