# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_LimitLOCKS_-1.target
[Unit]
Description=Foo
[Socket]
LimitLOCKS=-1

