# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ExecStopPost_a%3Eb.target
[Unit]
Description=Foo
[Socket]
ExecStopPost=a>b

