# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SendBuffer_-1T.target
[Unit]
Description=Foo
[Socket]
SendBuffer=-1T

