# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_FinalKillSignal_True.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=True

