# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_WorkingDirectory_ftp%3A----bar.gov.target
[Unit]
Description=Foo
[Socket]
WorkingDirectory=ftp://bar.gov

