# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SocketUser_meisnotausercozofmyloooooooooooooooooooooooooooooooooognamegee.target
[Unit]
Description=Foo
[Socket]
SocketUser=meisnotausercozofmyloooooooooooooooooooooooooooooooooognamegee

