# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_UtmpMode_1.target
[Unit]
Description=Foo
[Socket]
UtmpMode=1

