# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_User_1.target
[Unit]
Description=Foo
[Socket]
User=1

