# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_UMask_-1.target
[Unit]
Description=Foo
[Socket]
UMask=-1

