# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ExecStartPre_--etc--passwd.target
[Unit]
Description=Foo
[Socket]
ExecStartPre=/etc/passwd

