# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ExecStopPre_--etc--passwd.target
[Unit]
Description=Foo
[Socket]
ExecStopPre=/etc/passwd

