# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SendSIGKILL_bar.target
[Unit]
Description=Foo
[Socket]
SendSIGKILL=bar

