# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ReceiveBuffer_10000000000000K.target
[Unit]
Description=Foo
[Socket]
ReceiveBuffer=10000000000000K

