# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_BlockIODeviceWeight_--devvvv%200.target
[Unit]
Description=Foo
[Socket]
BlockIODeviceWeight=/devvvv 0

