# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SyslogFacility_True.target
[Unit]
Description=Foo
[Socket]
SyslogFacility=True

