# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ListenSpecial_True.target
[Unit]
Description=Foo
[Socket]
ListenSpecial=True

