# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_DeviceAllow_--devvvv.target
[Unit]
Description=Foo
[Socket]
DeviceAllow=/devvvv

