# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_KillMode_True.target
[Unit]
Description=Foo
[Socket]
KillMode=True

