# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_TimerSlackNSec_10000000000000.target
[Unit]
Description=Foo
[Socket]
TimerSlackNSec=10000000000000

