# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_KeepAliveProbes_10000000000000.target
[Unit]
Description=Foo
[Socket]
KeepAliveProbes=10000000000000

