# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_LimitNOFILE_10000000000000.target
[Unit]
Description=Foo
[Socket]
LimitNOFILE=10000000000000

