# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_KillMode_1.target
[Unit]
Description=Foo
[Socket]
KillMode=1

