# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_LimitCORE_10000000000000.target
[Unit]
Description=Foo
[Socket]
LimitCORE=10000000000000

