# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ExecStartPre_a%3Eb.target
[Unit]
Description=Foo
[Socket]
ExecStartPre=a>b

