# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SyslogLevel_1.target
[Unit]
Description=Foo
[Socket]
SyslogLevel=1

