# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_LimitLOCKS_10000000000000.target
[Unit]
Description=Foo
[Socket]
LimitLOCKS=10000000000000

