# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_IPTOS_-1.target
[Unit]
Description=Foo
[Socket]
IPTOS=-1

