# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_BlockIODeviceWeight_dev--foo%200.target
[Unit]
Description=Foo
[Socket]
BlockIODeviceWeight=dev/foo 0

