# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_DirectoryMode_999.target
[Unit]
Description=Foo
[Socket]
DirectoryMode=999

