# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_MountAPIVFS_baz.target
[Unit]
Description=Foo
[Socket]
MountAPIVFS=baz

