# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_DirectoryMode_True.target
[Unit]
Description=Foo
[Socket]
DirectoryMode=True

