# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_Symlinks_..--..--.target
[Unit]
Description=Foo
[Socket]
Symlinks=../../

