# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_TasksMax_0.target
[Unit]
Description=Foo
[Socket]
TasksMax=0

