# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_Personality_1.target
[Unit]
Description=Foo
[Socket]
Personality=1

