# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_MountFlags_1.target
[Unit]
Description=Foo
[Socket]
MountFlags=1

