# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ReceiveBuffer_-1B.target
[Unit]
Description=Foo
[Socket]
ReceiveBuffer=-1B

