# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ExecStartPost_True.target
[Unit]
Description=Foo
[Socket]
ExecStartPost=True

