# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_MountFlags_True.target
[Unit]
Description=Foo
[Socket]
MountFlags=True

