# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_MemoryDenyWriteExecute_bar.target
[Unit]
Description=Foo
[Socket]
MemoryDenyWriteExecute=bar

