# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_RootImage_ftp%3A----bar.gov.target
[Unit]
Description=Foo
[Socket]
RootImage=ftp://bar.gov

