# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_KeepAliveIntervalSec_10000000000000.target
[Unit]
Description=Foo
[Socket]
KeepAliveIntervalSec=10000000000000

