# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_User_True.target
[Unit]
Description=Foo
[Socket]
User=True

