# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SupplementaryGroups_1.target
[Unit]
Description=Foo
[Socket]
SupplementaryGroups=1

