# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_Broadcast_baz.target
[Unit]
Description=Foo
[Socket]
Broadcast=baz

