# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SocketProtocol_True.target
[Unit]
Description=Foo
[Socket]
SocketProtocol=True

