# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ExecStopPre_True.target
[Unit]
Description=Foo
[Socket]
ExecStopPre=True

