# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_User_meisnotausercozofmyloooooooooooooooooooooooooooooooooognamegee.target
[Unit]
Description=Foo
[Socket]
User=meisnotausercozofmyloooooooooooooooooooooooooooooooooognamegee

