# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_Transparent_bar.target
[Unit]
Description=Foo
[Socket]
Transparent=bar

