# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_RemoveOnStop_foo.target
[Unit]
Description=Foo
[Socket]
RemoveOnStop=foo

