# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_BlockIOAccounting_baz.target
[Unit]
Description=Foo
[Socket]
BlockIOAccounting=baz

