# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_Personality_True.target
[Unit]
Description=Foo
[Socket]
Personality=True

