# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ReusePort_foo.target
[Unit]
Description=Foo
[Socket]
ReusePort=foo

