# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ExecStopPost_1.target
[Unit]
Description=Foo
[Socket]
ExecStopPost=1

