# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_Mark_-1.target
[Unit]
Description=Foo
[Socket]
Mark=-1

