# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SELinuxContextFromNet_bar.target
[Unit]
Description=Foo
[Socket]
SELinuxContextFromNet=bar

