# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ExecStartPre_True.target
[Unit]
Description=Foo
[Socket]
ExecStartPre=True

