# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_PrivateNetwork_foo.target
[Unit]
Description=Foo
[Socket]
PrivateNetwork=foo

