# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_Writable_foo.target
[Unit]
Description=Foo
[Socket]
Writable=foo

