# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_Backlog_-1.target
[Unit]
Description=Foo
[Socket]
Backlog=-1

