# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_RemoveIPC_foo.target
[Unit]
Description=Foo
[Socket]
RemoveIPC=foo

