# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SocketGroup_meisnotagroupcozofmyloooooooooooooooooooooooooooooooooognamegee.target
[Unit]
Description=Foo
[Socket]
SocketGroup=meisnotagroupcozofmyloooooooooooooooooooooooooooooooooognamegee

