# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SyslogLevelPrefix_foo.target
[Unit]
Description=Foo
[Socket]
SyslogLevelPrefix=foo

