# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_ExecStopPre_1.target
[Unit]
Description=Foo
[Socket]
ExecStopPre=1

