# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_SocketUser_1.target
[Unit]
Description=Foo
[Socket]
SocketUser=1

