# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_RemoveIPC_bar.target
[Unit]
Description=Foo
[Socket]
RemoveIPC=bar

