# systemdlint/tests/auto/InvalidValue/2.31/bad/Socket_DirectoryMode_1.target
[Unit]
Description=Foo
[Socket]
DirectoryMode=1

