# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_Delegate_True.target
[Unit]
Description=Foo
[Swap]
Delegate=True

