# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_LimitMSGQUEUE_10000000000000.target
[Unit]
Description=Foo
[Swap]
LimitMSGQUEUE=10000000000000

