# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_LimitFSIZE_-1.target
[Unit]
Description=Foo
[Swap]
LimitFSIZE=-1

