# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_UMask_-1.target
[Unit]
Description=Foo
[Swap]
UMask=-1

