# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_Priority_-1.target
[Unit]
Description=Foo
[Swap]
Priority=-1

