# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_KillSignal_True.target
[Unit]
Description=Foo
[Swap]
KillSignal=True

