# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_DisableControllers_1.target
[Unit]
Description=Foo
[Swap]
DisableControllers=1

