# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_RootImage_http%3A----foo.org.target
[Unit]
Description=Foo
[Swap]
RootImage=http://foo.org

