# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_RootDirectory_ftp%3A----bar.gov.target
[Unit]
Description=Foo
[Swap]
RootDirectory=ftp://bar.gov

