# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_TTYPath_--tmp--dev.target
[Unit]
Description=Foo
[Swap]
TTYPath=/tmp/dev

