# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_WatchdogSignal_True.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=True

