# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_Group_True.target
[Unit]
Description=Foo
[Swap]
Group=True

