# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_CPUSchedulingPolicy_1.target
[Unit]
Description=Foo
[Swap]
CPUSchedulingPolicy=1

