# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_LimitLOCKS_10000000000000.target
[Unit]
Description=Foo
[Swap]
LimitLOCKS=10000000000000

