# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_SyslogLevelPrefix_bar.target
[Unit]
Description=Foo
[Swap]
SyslogLevelPrefix=bar

