# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_PrivateNetwork_bar.target
[Unit]
Description=Foo
[Swap]
PrivateNetwork=bar

