# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_Environment_True.target
[Unit]
Description=Foo
[Swap]
Environment=True

