# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_IgnoreSIGPIPE_baz.target
[Unit]
Description=Foo
[Swap]
IgnoreSIGPIPE=baz

