# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_BlockIOWeight_10001.target
[Unit]
Description=Foo
[Swap]
BlockIOWeight=10001

