# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_IOWeight_0.target
[Unit]
Description=Foo
[Swap]
IOWeight=0

