# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_IOSchedulingClass_True.target
[Unit]
Description=Foo
[Swap]
IOSchedulingClass=True

