# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_SecureBits_1.target
[Unit]
Description=Foo
[Swap]
SecureBits=1

