# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_TTYVHangup_baz.target
[Unit]
Description=Foo
[Swap]
TTYVHangup=baz

