# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_BlockIODeviceWeight_--tmp--dev%2010001.target
[Unit]
Description=Foo
[Swap]
BlockIODeviceWeight=/tmp/dev 10001

