# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_BlockIOReadBandwidthMax_--devvvv%20-1K.target
[Unit]
Description=Foo
[Swap]
BlockIOReadBandwidthMax=/devvvv -1K

