# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_Environment_1.target
[Unit]
Description=Foo
[Swap]
Environment=1

