# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_TTYReset_baz.target
[Unit]
Description=Foo
[Swap]
TTYReset=baz

