# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_DisableControllers_True.target
[Unit]
Description=Foo
[Swap]
DisableControllers=True

