# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_WorkingDirectory_http%3A----foo.org.target
[Unit]
Description=Foo
[Swap]
WorkingDirectory=http://foo.org

