# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_TTYReset_foo.target
[Unit]
Description=Foo
[Swap]
TTYReset=foo

