# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_FinalKillSignal_True.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=True

