# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_MemoryAccounting_foo.target
[Unit]
Description=Foo
[Swap]
MemoryAccounting=foo

