# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_TasksAccounting_baz.target
[Unit]
Description=Foo
[Swap]
TasksAccounting=baz

