# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_TTYVTDisallocate_baz.target
[Unit]
Description=Foo
[Swap]
TTYVTDisallocate=baz

