# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_IOSchedulingClass_1.target
[Unit]
Description=Foo
[Swap]
IOSchedulingClass=1

