# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_MemoryHigh_-1M.target
[Unit]
Description=Foo
[Swap]
MemoryHigh=-1M

