# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_CPUSchedulingResetOnFork_baz.target
[Unit]
Description=Foo
[Swap]
CPUSchedulingResetOnFork=baz

