# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_MemorySwapMax_-1T.target
[Unit]
Description=Foo
[Swap]
MemorySwapMax=-1T

