# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_CPUSchedulingResetOnFork_foo.target
[Unit]
Description=Foo
[Swap]
CPUSchedulingResetOnFork=foo

