# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_LimitSTACK_-1.target
[Unit]
Description=Foo
[Swap]
LimitSTACK=-1

