# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_RootDirectory_http%3A----foo.org.target
[Unit]
Description=Foo
[Swap]
RootDirectory=http://foo.org

