# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_SyslogFacility_True.target
[Unit]
Description=Foo
[Swap]
SyslogFacility=True

