# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_LimitNPROC_-1.target
[Unit]
Description=Foo
[Swap]
LimitNPROC=-1

