# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_Delegate_1.target
[Unit]
Description=Foo
[Swap]
Delegate=1

