# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_MountAPIVFS_bar.target
[Unit]
Description=Foo
[Swap]
MountAPIVFS=bar

