# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_BlockIOAccounting_baz.target
[Unit]
Description=Foo
[Swap]
BlockIOAccounting=baz

