# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_Priority_10000000000000.target
[Unit]
Description=Foo
[Swap]
Priority=10000000000000

