# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_User_meisnotausercozofmyloooooooooooooooooooooooooooooooooognamegee.target
[Unit]
Description=Foo
[Swap]
User=meisnotausercozofmyloooooooooooooooooooooooooooooooooognamegee

