# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_UtmpMode_True.target
[Unit]
Description=Foo
[Swap]
UtmpMode=True

