# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_RootImage_ftp%3A----bar.gov.target
[Unit]
Description=Foo
[Swap]
RootImage=ftp://bar.gov

