# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_KillMode_1.target
[Unit]
Description=Foo
[Swap]
KillMode=1

