# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_BlockIOWriteBandwidthMax_--tmp--dev%20-1B.target
[Unit]
Description=Foo
[Swap]
BlockIOWriteBandwidthMax=/tmp/dev -1B

