# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_MemorySwapMax_-1B.target
[Unit]
Description=Foo
[Swap]
MemorySwapMax=-1B

