# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_KillMode_True.target
[Unit]
Description=Foo
[Swap]
KillMode=True

