# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_BlockIODeviceWeight_--devvvv%200.target
[Unit]
Description=Foo
[Swap]
BlockIODeviceWeight=/devvvv 0

