# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_SendSIGKILL_foo.target
[Unit]
Description=Foo
[Swap]
SendSIGKILL=foo

