# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_RemoveIPC_foo.target
[Unit]
Description=Foo
[Swap]
RemoveIPC=foo

