# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_TimeoutSec_1p.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=1p

