# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_IOAccounting_bar.target
[Unit]
Description=Foo
[Swap]
IOAccounting=bar

