# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_BlockIOWriteBandwidthMax_--devvvv%2010000000000000B.target
[Unit]
Description=Foo
[Swap]
BlockIOWriteBandwidthMax=/devvvv 10000000000000B

