# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_SecureBits_True.target
[Unit]
Description=Foo
[Swap]
SecureBits=True

