# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_LimitCORE_-1.target
[Unit]
Description=Foo
[Swap]
LimitCORE=-1

