# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_RemoveIPC_bar.target
[Unit]
Description=Foo
[Swap]
RemoveIPC=bar

