# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_IOSchedulingPriority_8.target
[Unit]
Description=Foo
[Swap]
IOSchedulingPriority=8

