# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_SupplementaryGroups_1.target
[Unit]
Description=Foo
[Swap]
SupplementaryGroups=1

