# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_MemorySwapMax_-1K.target
[Unit]
Description=Foo
[Swap]
MemorySwapMax=-1K

