# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_BlockIOWeight_0.target
[Unit]
Description=Foo
[Swap]
BlockIOWeight=0

