# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_BlockIOWriteBandwidthMax_--devvvv%20-1G.target
[Unit]
Description=Foo
[Swap]
BlockIOWriteBandwidthMax=/devvvv -1G

