# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_TasksMax_0.target
[Unit]
Description=Foo
[Swap]
TasksMax=0

