# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_LimitCORE_10000000000000.target
[Unit]
Description=Foo
[Swap]
LimitCORE=10000000000000

