# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_MemoryLow_10000000000000B.target
[Unit]
Description=Foo
[Swap]
MemoryLow=10000000000000B

