# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_WorkingDirectory_ftp%3A----bar.gov.target
[Unit]
Description=Foo
[Swap]
WorkingDirectory=ftp://bar.gov

