# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_MountFlags_True.target
[Unit]
Description=Foo
[Swap]
MountFlags=True

