# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_EnvironmentFile_True.target
[Unit]
Description=Foo
[Swap]
EnvironmentFile=True

