# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_SyslogLevel_1.target
[Unit]
Description=Foo
[Swap]
SyslogLevel=1

