# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_Personality_1.target
[Unit]
Description=Foo
[Swap]
Personality=1

