# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_LimitCPU_10000000000000.target
[Unit]
Description=Foo
[Swap]
LimitCPU=10000000000000

