# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_DeviceAllow_--tmp.target
[Unit]
Description=Foo
[Swap]
DeviceAllow=/tmp

