# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_BlockIOWriteBandwidthMax_--tmp%20-1M.target
[Unit]
Description=Foo
[Swap]
BlockIOWriteBandwidthMax=/tmp -1M

