# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_DeviceAllow_--devvvv.target
[Unit]
Description=Foo
[Swap]
DeviceAllow=/devvvv

