# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_TimerSlackNSec_10000000000000.target
[Unit]
Description=Foo
[Swap]
TimerSlackNSec=10000000000000

