# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_SystemCallErrorNumber_4096.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=4096

