# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_BlockIOWriteBandwidthMax_--tmp--dev%2010000000000000G.target
[Unit]
Description=Foo
[Swap]
BlockIOWriteBandwidthMax=/tmp/dev 10000000000000G

