# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_SyslogFacility_1.target
[Unit]
Description=Foo
[Swap]
SyslogFacility=1

