# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_SendSIGHUP_baz.target
[Unit]
Description=Foo
[Swap]
SendSIGHUP=baz

