# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_DevicePolicy_1.target
[Unit]
Description=Foo
[Swap]
DevicePolicy=1

