# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_PrivateDevices_foo.target
[Unit]
Description=Foo
[Swap]
PrivateDevices=foo

