# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_LimitRTTIME_10000000000000.target
[Unit]
Description=Foo
[Swap]
LimitRTTIME=10000000000000

