# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_DevicePolicy_True.target
[Unit]
Description=Foo
[Swap]
DevicePolicy=True

