# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_MemoryDenyWriteExecute_foo.target
[Unit]
Description=Foo
[Swap]
MemoryDenyWriteExecute=foo

