# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_LimitNICE_10000000000000.target
[Unit]
Description=Foo
[Swap]
LimitNICE=10000000000000

