# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_Personality_True.target
[Unit]
Description=Foo
[Swap]
Personality=True

