# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_PrivateTmp_foo.target
[Unit]
Description=Foo
[Swap]
PrivateTmp=foo

