# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_UtmpMode_1.target
[Unit]
Description=Foo
[Swap]
UtmpMode=1

