# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_CPUSchedulingPolicy_True.target
[Unit]
Description=Foo
[Swap]
CPUSchedulingPolicy=True

