# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_MountFlags_1.target
[Unit]
Description=Foo
[Swap]
MountFlags=1

