# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_BlockIOReadBandwidthMax_--tmp--dev%20-1G.target
[Unit]
Description=Foo
[Swap]
BlockIOReadBandwidthMax=/tmp/dev -1G

