# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_IOAccounting_baz.target
[Unit]
Description=Foo
[Swap]
IOAccounting=baz

