# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_StartupCPUShares_262145.target
[Unit]
Description=Foo
[Swap]
StartupCPUShares=262145

