# systemdlint/tests/auto/InvalidValue/2.31/bad/Swap_MemoryMax_-1B.target
[Unit]
Description=Foo
[Swap]
MemoryMax=-1B

