# systemdlint/tests/auto/InvalidValue/2.31/bad/Tap_MultiQueue_bar.target
[Unit]
Description=Foo
[Tap]
MultiQueue=bar

