# systemdlint/tests/auto/InvalidValue/2.31/bad/Tap_OneQueue_foo.target
[Unit]
Description=Foo
[Tap]
OneQueue=foo

