# systemdlint/tests/auto/InvalidValue/2.31/bad/Tap_MultiQueue_foo.target
[Unit]
Description=Foo
[Tap]
MultiQueue=foo

