# systemdlint/tests/auto/InvalidValue/2.31/bad/Tap_OneQueue_bar.target
[Unit]
Description=Foo
[Tap]
OneQueue=bar

