# systemdlint/tests/auto/InvalidValue/2.31/bad/Timer_OnStartupSec_3decades.target
[Unit]
Description=Foo
[Timer]
OnStartupSec=3decades

