# systemdlint/tests/auto/InvalidValue/2.31/bad/Timer_AccuracySec_2potatoes.target
[Unit]
Description=Foo
[Timer]
AccuracySec=2potatoes

