# systemdlint/tests/auto/InvalidValue/2.31/bad/Timer_RemainAfterElapse_baz.target
[Unit]
Description=Foo
[Timer]
RemainAfterElapse=baz

