# systemdlint/tests/auto/InvalidValue/2.31/bad/Timer_OnActiveSec_3decades.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=3decades

