# systemdlint/tests/auto/InvalidValue/2.31/bad/Timer_Persistent_baz.target
[Unit]
Description=Foo
[Timer]
Persistent=baz

