# systemdlint/tests/auto/InvalidValue/2.31/bad/Timer_AccuracySec_True.target
[Unit]
Description=Foo
[Timer]
AccuracySec=True

