# systemdlint/tests/auto/InvalidValue/2.31/bad/Timer_OnActiveSec_1p.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=1p

