# systemdlint/tests/auto/InvalidValue/2.31/bad/Timer_OnClockChange_bar.target
[Unit]
Description=Foo
[Timer]
OnClockChange=bar

