# systemdlint/tests/auto/InvalidValue/2.31/bad/Timer_OnUnitActiveSec_True.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=True

