# systemdlint/tests/auto/InvalidValue/2.31/bad/Time_PollIntervalMaxSec_10000000000000.target
[Unit]
Description=Foo
[Time]
PollIntervalMaxSec=10000000000000

