# systemdlint/tests/auto/InvalidValue/2.31/bad/Time_RootDistanceMaxSec_-1.target
[Unit]
Description=Foo
[Time]
RootDistanceMaxSec=-1

