# systemdlint/tests/auto/InvalidValue/2.31/bad/Timer_OnBootSec_True.target
[Unit]
Description=Foo
[Timer]
OnBootSec=True

