# systemdlint/tests/auto/InvalidValue/2.31/bad/Timer_RandomizedDelaySec_2potatoes.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=2potatoes

