# systemdlint/tests/auto/InvalidValue/2.31/bad/Time_RootDistanceMaxSec_10000000000000.target
[Unit]
Description=Foo
[Time]
RootDistanceMaxSec=10000000000000

