# systemdlint/tests/auto/InvalidValue/2.31/bad/Time_PollIntervalMinSec_15.target
[Unit]
Description=Foo
[Time]
PollIntervalMinSec=15

