# systemdlint/tests/auto/InvalidValue/2.31/bad/Timer_RemainAfterElapse_bar.target
[Unit]
Description=Foo
[Timer]
RemainAfterElapse=bar

