# systemdlint/tests/auto/InvalidValue/2.31/bad/Timer_OnStartupSec_1p.target
[Unit]
Description=Foo
[Timer]
OnStartupSec=1p

