# systemdlint/tests/auto/InvalidValue/2.31/bad/Timer_OnTimezoneChange_foo.target
[Unit]
Description=Foo
[Timer]
OnTimezoneChange=foo

