# systemdlint/tests/auto/InvalidValue/2.31/bad/Timer_RemainAfterElapse_foo.target
[Unit]
Description=Foo
[Timer]
RemainAfterElapse=foo

