# systemdlint/tests/auto/InvalidValue/2.31/bad/Tun_OneQueue_bar.target
[Unit]
Description=Foo
[Tun]
OneQueue=bar

