# systemdlint/tests/auto/InvalidValue/2.31/bad/Tun_MultiQueue_foo.target
[Unit]
Description=Foo
[Tun]
MultiQueue=foo

