# systemdlint/tests/auto/InvalidValue/2.31/bad/Tun_Group_meisnotagroupcozofmyloooooooooooooooooooooooooooooooooognamegee.target
[Unit]
Description=Foo
[Tun]
Group=meisnotagroupcozofmyloooooooooooooooooooooooooooooooooognamegee

