# systemdlint/tests/auto/InvalidValue/2.31/bad/Tun_PacketInfo_foo.target
[Unit]
Description=Foo
[Tun]
PacketInfo=foo

