# systemdlint/tests/auto/InvalidValue/2.31/bad/Tun_VNetHeader_foo.target
[Unit]
Description=Foo
[Tun]
VNetHeader=foo

