# systemdlint/tests/auto/InvalidValue/2.31/bad/Tun_VNetHeader_bar.target
[Unit]
Description=Foo
[Tun]
VNetHeader=bar

