# systemdlint/tests/auto/InvalidValue/2.31/bad/Tun_PacketInfo_bar.target
[Unit]
Description=Foo
[Tun]
PacketInfo=bar

