# systemdlint/tests/auto/InvalidValue/2.31/bad/Tun_User_True.target
[Unit]
Description=Foo
[Tun]
User=True

