# systemdlint/tests/auto/InvalidValue/2.31/bad/Tunnel_AllowLocalRemote_foo.target
[Unit]
Description=Foo
[Tunnel]
AllowLocalRemote=foo

