# systemdlint/tests/auto/InvalidValue/2.31/bad/Tunnel_Independent_baz.target
[Unit]
Description=Foo
[Tunnel]
Independent=baz

