# systemdlint/tests/auto/InvalidValue/2.31/bad/Tunnel_AllowLocalRemote_baz.target
[Unit]
Description=Foo
[Tunnel]
AllowLocalRemote=baz

