# systemdlint/tests/auto/InvalidValue/2.31/bad/Tunnel_Mode_1.target
[Unit]
Description=Foo
[Tunnel]
Mode=1

