# systemdlint/tests/auto/InvalidValue/2.31/bad/Tunnel_FooOverUDP_bar.target
[Unit]
Description=Foo
[Tunnel]
FooOverUDP=bar

