# systemdlint/tests/auto/InvalidValue/2.31/bad/Tunnel_ERSPANIndex_1048576.target
[Unit]
Description=Foo
[Tunnel]
ERSPANIndex=1048576

