# systemdlint/tests/auto/InvalidValue/2.31/bad/Tunnel_CopyDSCP_bar.target
[Unit]
Description=Foo
[Tunnel]
CopyDSCP=bar

