# systemdlint/tests/auto/InvalidValue/2.31/bad/Tunnel_Mode_True.target
[Unit]
Description=Foo
[Tunnel]
Mode=True

