# systemdlint/tests/auto/InvalidValue/2.31/bad/Tunnel_SerializeTunneledPackets_baz.target
[Unit]
Description=Foo
[Tunnel]
SerializeTunneledPackets=baz

