# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_AssertACPower_foo.target
[Unit]
Description=Foo
Description=abc
AssertACPower=foo

