# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_ConditionDirectoryNotEmpty_1.target
[Unit]
Description=Foo
Description=abc
ConditionDirectoryNotEmpty=1

