# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_ConditionPathExists_1.target
[Unit]
Description=Foo
Description=abc
ConditionPathExists=1

