# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_JobTimeoutAction_1.target
[Unit]
Description=Foo
Description=abc
JobTimeoutAction=1

