# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_FailureAction_1.target
[Unit]
Description=Foo
Description=abc
FailureAction=1

