# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_JobTimeoutRebootArgument_1.target
[Unit]
Description=Foo
Description=abc
JobTimeoutRebootArgument=1

