# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_AssertPathIsReadWrite_1.target
[Unit]
Description=Foo
Description=abc
AssertPathIsReadWrite=1

