# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_OnFailureJobMode_1.target
[Unit]
Description=Foo
Description=abc
OnFailureJobMode=1

