# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_ConditionFileNotEmpty_tmp.target
[Unit]
Description=Foo
Description=abc
ConditionFileNotEmpty=tmp

