# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_RefuseManualStart_baz.target
[Unit]
Description=Foo
Description=abc
RefuseManualStart=baz

