# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_ConditionACPower_baz.target
[Unit]
Description=Foo
Description=abc
ConditionACPower=baz

