# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_ConditionPathIsMountPoint_True.target
[Unit]
Description=Foo
Description=abc
ConditionPathIsMountPoint=True

