# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_ConditionPathIsDirectory_..--..--.target
[Unit]
Description=Foo
Description=abc
ConditionPathIsDirectory=../../

