# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_RequiresMountsFor_1.target
[Unit]
Description=Foo
Description=abc
RequiresMountsFor=1

