# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_AssertPathExists_True.target
[Unit]
Description=Foo
Description=abc
AssertPathExists=True

