# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_IgnoreOnIsolate_baz.target
[Unit]
Description=Foo
Description=abc
IgnoreOnIsolate=baz

