# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_AssertVirtualization_bar.target
[Unit]
Description=Foo
Description=abc
AssertVirtualization=bar

