# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_ConditionPathExists_True.target
[Unit]
Description=Foo
Description=abc
ConditionPathExists=True

