# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_ConditionDirectoryNotEmpty_True.target
[Unit]
Description=Foo
Description=abc
ConditionDirectoryNotEmpty=True

