# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_AssertCPUs_10000000000000.target
[Unit]
Description=Foo
Description=abc
AssertCPUs=10000000000000

