# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_AssertPathIsDirectory_1.target
[Unit]
Description=Foo
Description=abc
AssertPathIsDirectory=1

