# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_ConditionPathIsSymbolicLink_..--..--.target
[Unit]
Description=Foo
Description=abc
ConditionPathIsSymbolicLink=../../

