# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_AssertPathIsSymbolicLink_1.target
[Unit]
Description=Foo
Description=abc
AssertPathIsSymbolicLink=1

