# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_StartLimitAction_1.target
[Unit]
Description=Foo
Description=abc
StartLimitAction=1

