# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_JobRunningTimeoutSec_1p.target
[Unit]
Description=Foo
Description=abc
JobRunningTimeoutSec=1p

