# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_AssertFileNotEmpty_True.target
[Unit]
Description=Foo
Description=abc
AssertFileNotEmpty=True

