# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_RequiresMountsFor_tmp.target
[Unit]
Description=Foo
Description=abc
RequiresMountsFor=tmp

