# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_AssertACPower_bar.target
[Unit]
Description=Foo
Description=abc
AssertACPower=bar

