# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_AssertDirectoryNotEmpty_..--..--.target
[Unit]
Description=Foo
Description=abc
AssertDirectoryNotEmpty=../../

