# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_StartLimitAction_True.target
[Unit]
Description=Foo
Description=abc
StartLimitAction=True

