# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_ConditionFileNotEmpty_1.target
[Unit]
Description=Foo
Description=abc
ConditionFileNotEmpty=1

