# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_JobTimeoutRebootArgument_True.target
[Unit]
Description=Foo
Description=abc
JobTimeoutRebootArgument=True

