# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_ConditionPathIsMountPoint_tmp.target
[Unit]
Description=Foo
Description=abc
ConditionPathIsMountPoint=tmp

