# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_JobTimeoutAction_True.target
[Unit]
Description=Foo
Description=abc
JobTimeoutAction=True

