# systemdlint/tests/auto/InvalidValue/2.31/bad/Unit_ConditionPathIsSymbolicLink_1.target
[Unit]
Description=Foo
Description=abc
ConditionPathIsSymbolicLink=1

