# systemdlint/tests/auto/InvalidValue/2.31/bad/VLAN_MVRP_bar.target
[Unit]
Description=Foo
[VLAN]
MVRP=bar

