# systemdlint/tests/auto/InvalidValue/2.31/bad/VLAN_LooseBinding_baz.target
[Unit]
Description=Foo
[VLAN]
LooseBinding=baz

