# systemdlint/tests/auto/InvalidValue/2.31/bad/VLAN_GVRP_bar.target
[Unit]
Description=Foo
[VLAN]
GVRP=bar

