# systemdlint/tests/auto/InvalidValue/2.31/bad/VLAN_GVRP_baz.target
[Unit]
Description=Foo
[VLAN]
GVRP=baz

