# systemdlint/tests/auto/InvalidValue/2.31/bad/VLAN_LooseBinding_foo.target
[Unit]
Description=Foo
[VLAN]
LooseBinding=foo

