# systemdlint/tests/auto/InvalidValue/2.31/bad/VXLAN_ReduceARPProxy_baz.target
[Unit]
Description=Foo
[VXLAN]
ReduceARPProxy=baz

