# systemdlint/tests/auto/InvalidValue/2.31/bad/VXLAN_Remote_0--abc.target
[Unit]
Description=Foo
[VXLAN]
Remote=0/abc

