# systemdlint/tests/auto/InvalidValue/2.31/bad/VXLAN_FDBAgeingSec_10000000000000.target
[Unit]
Description=Foo
[VXLAN]
FDBAgeingSec=10000000000000

