# systemdlint/tests/auto/InvalidValue/2.31/bad/VXLAN_RemoteChecksumRx_baz.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumRx=baz

