# systemdlint/tests/auto/InvalidValue/2.31/bad/VXLAN_MacLearning_bar.target
[Unit]
Description=Foo
[VXLAN]
MacLearning=bar

