# systemdlint/tests/auto/InvalidValue/2.31/bad/VXLAN_Remote_abc.target
[Unit]
Description=Foo
[VXLAN]
Remote=abc

