# systemdlint/tests/auto/InvalidValue/2.31/bad/VXLAN_L3MissNotification_bar.target
[Unit]
Description=Foo
[VXLAN]
L3MissNotification=bar

