# systemdlint/tests/auto/InvalidValue/2.31/bad/VXLAN_UDP6ZeroChecksumRx_foo.target
[Unit]
Description=Foo
[VXLAN]
UDP6ZeroChecksumRx=foo

