# systemdlint/tests/auto/InvalidValue/2.31/bad/VXLAN_Local_abc.target
[Unit]
Description=Foo
[VXLAN]
Local=abc

