# systemdlint/tests/auto/InvalidValue/2.31/bad/VXLAN_Remote_255.256.255.255.target
[Unit]
Description=Foo
[VXLAN]
Remote=255.256.255.255

