# systemdlint/tests/auto/InvalidValue/2.31/bad/WireGuardPeer_PresharedKeyFile_1.target
[Unit]
Description=Foo
[WireGuardPeer]
PresharedKeyFile=1

