# systemdlint/tests/auto/InvalidValue/2.31/bad/WireGuardPeer_PersistentKeepalive_65536.target
[Unit]
Description=Foo
[WireGuardPeer]
PersistentKeepalive=65536

