# tests/auto/OptionDeprecated/2.46/bad/TrafficControlQueueingDiscipline_NetworkEmulatorDelaySec_0.target
[Unit]
Description=Foo
[TrafficControlQueueingDiscipline]
NetworkEmulatorDelaySec=123

