# tests/auto/OptionDeprecated/2.46/bad/Unit_LogNamespace_abc.target
[Unit]
Description=Foo
Description=abc
LogNamespace=123

