# tests/auto/InvalidValue/2.28/good/Unit_RequiresOverridable_abc.target
[Unit]
Description=Foo
Description=abc
RequiresOverridable=abc

