# systemdlint/tests/auto/OptionTooNew/2.14/good/Service_ProtectHome_tmpfs.target
[Unit]
Description=Foo
[Service]
ProtectHome=123

