# systemdlint/tests/auto/OptionTooNew/2.14/good/Socket_Symlinks_--tmp.target
[Unit]
Description=Foo
[Socket]
Symlinks=123

