# systemdlint/tests/auto/OptionTooNew/2.14/good/Socket_SocketUser_iamauser.target
[Unit]
Description=Foo
[Socket]
SocketUser=123

