# systemdlint/tests/auto/OptionTooNew/2.14/good/Swap_ProtectSystem_full.target
[Unit]
Description=Foo
[Swap]
ProtectSystem=123

