# systemdlint/tests/auto/OptionTooNew/2.14/good/Swap_ProtectHome_tmpfs.target
[Unit]
Description=Foo
[Swap]
ProtectHome=123

