# systemdlint/tests/auto/OptionDeprecated/2.31/bad/Swap_Capabilities_abc.target
[Unit]
Description=Foo
[Swap]
Capabilities=123

