# tests/auto/InvalidValue/2.12/good/Mount_ControlGroup_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
ControlGroup=abc

