# tests/auto/InvalidValue/2.12/good/Service_ControlGroupAttribute_abc.target
[Unit]
Description=Foo
[Service]
ControlGroupAttribute=abc

