# tests/auto/InvalidValue/2.12/good/Service_ControlGroupModify_abc.target
[Unit]
Description=Foo
[Service]
ControlGroupModify=abc

