# tests/auto/InvalidValue/2.12/good/Service_ControlGroupPersistent_abc.target
[Unit]
Description=Foo
[Service]
ControlGroupPersistent=abc

