# tests/auto/InvalidValue/2.12/good/Socket_ControlGroupModify_abc.target
[Unit]
Description=Foo
[Socket]
ControlGroupModify=abc

