# tests/auto/InvalidValue/2.12/good/Socket_ControlGroup_abc.target
[Unit]
Description=Foo
[Socket]
ControlGroup=abc

