# tests/auto/InvalidValue/2.12/good/Socket_ControlGroupPersistent_abc.target
[Unit]
Description=Foo
[Socket]
ControlGroupPersistent=abc

