# tests/auto/InvalidValue/2.12/good/Swap_ControlGroupModify_abc.target
[Unit]
Description=Foo
[Swap]
ControlGroupModify=abc

