# tests/auto/InvalidValue/2.31/bad/Address_Broadcast_1.1.1.1---1.target
[Unit]
Description=Foo
[Address]
Address=abc
Broadcast=1.1.1.1/-1

