# tests/auto/InvalidValue/2.31/bad/Address_Peer_0--abc.target
[Unit]
Description=Foo
[Address]
Address=abc
Peer=0/abc

