# tests/auto/InvalidValue/2.31/bad/Address_DuplicateAddressDetection_baz.target
[Unit]
Description=Foo
[Address]
Address=abc
DuplicateAddressDetection=baz

