# tests/auto/InvalidValue/2.31/bad/Address_Peer_abc.target
[Unit]
Description=Foo
[Address]
Address=abc
Peer=abc

