# tests/auto/InvalidValue/2.31/bad/Address_Broadcast_255.256.255.255.target
[Unit]
Description=Foo
[Address]
Address=abc
Broadcast=255.256.255.255

