# tests/auto/InvalidValue/2.31/bad/Automount_TimeoutIdleSec_0.5d.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=0.5d

