# tests/auto/InvalidValue/2.31/bad/Automount_Where_1.target
[Unit]
Description=Foo
[Automount]
Where=1

