# tests/auto/InvalidValue/2.31/bad/Bond_UpDelaySec_-1.target
[Unit]
Description=Foo
[Bond]
UpDelaySec=-1

