# tests/auto/InvalidValue/2.31/bad/Bond_AdUserPortKey_0.target
[Unit]
Description=Foo
[Bond]
AdUserPortKey=0

