# tests/auto/InvalidValue/2.31/bad/Bond_GratuitousARP_256.target
[Unit]
Description=Foo
[Bond]
GratuitousARP=256

