# tests/auto/InvalidValue/2.31/bad/Bond_DynamicTransmitLoadBalancing_foo.target
[Unit]
Description=Foo
[Bond]
DynamicTransmitLoadBalancing=foo

