# tests/auto/InvalidValue/2.31/bad/Bond_GratuitousARP_-1.target
[Unit]
Description=Foo
[Bond]
GratuitousARP=-1

