# tests/auto/InvalidValue/2.31/bad/Bond_ARPIntervalSec_-1.target
[Unit]
Description=Foo
[Bond]
ARPIntervalSec=-1

