# tests/auto/InvalidValue/2.31/bad/Bond_MIIMonitorSec_10000000000000.target
[Unit]
Description=Foo
[Bond]
MIIMonitorSec=10000000000000

