# tests/auto/InvalidValue/2.31/bad/Bond_PacketsPerSlave_65536.target
[Unit]
Description=Foo
[Bond]
PacketsPerSlave=65536

