# tests/auto/InvalidValue/2.31/bad/Bond_AdActorSystemPriority_0.target
[Unit]
Description=Foo
[Bond]
AdActorSystemPriority=0

