# tests/auto/InvalidValue/2.31/bad/Bond_ResendIGMP_-1.target
[Unit]
Description=Foo
[Bond]
ResendIGMP=-1

