# tests/auto/InvalidValue/2.31/bad/Bond_LearnPacketIntervalSec_0.target
[Unit]
Description=Foo
[Bond]
LearnPacketIntervalSec=0

