# tests/auto/InvalidValue/2.31/bad/Bond_DynamicTransmitLoadBalancing_baz.target
[Unit]
Description=Foo
[Bond]
DynamicTransmitLoadBalancing=baz

