# tests/auto/InvalidValue/2.31/bad/Bond_UpDelaySec_10000000000000.target
[Unit]
Description=Foo
[Bond]
UpDelaySec=10000000000000

