# tests/auto/InvalidValue/2.31/bad/Bond_DynamicTransmitLoadBalancing_bar.target
[Unit]
Description=Foo
[Bond]
DynamicTransmitLoadBalancing=bar

