# tests/auto/InvalidValue/2.31/bad/Bond_DownDelaySec_10000000000000.target
[Unit]
Description=Foo
[Bond]
DownDelaySec=10000000000000

