# tests/auto/InvalidValue/2.31/bad/Bridge_Learning_foo.target
[Unit]
Description=Foo
[Bridge]
Learning=foo

