# tests/auto/InvalidValue/2.31/bad/Bridge_MulticastFlood_bar.target
[Unit]
Description=Foo
[Bridge]
MulticastFlood=bar

