# tests/auto/InvalidValue/2.31/bad/Bridge_Priority_-1.target
[Unit]
Description=Foo
[Bridge]
Priority=-1

