# tests/auto/InvalidValue/2.31/bad/Bridge_NeighborSuppression_bar.target
[Unit]
Description=Foo
[Bridge]
NeighborSuppression=bar

