# tests/auto/InvalidValue/2.31/bad/Bridge_MulticastFlood_baz.target
[Unit]
Description=Foo
[Bridge]
MulticastFlood=baz

