# tests/auto/InvalidValue/2.31/bad/Bridge_UnicastFlood_baz.target
[Unit]
Description=Foo
[Bridge]
UnicastFlood=baz

