# tests/auto/InvalidValue/2.31/bad/Bridge_UnicastFlood_foo.target
[Unit]
Description=Foo
[Bridge]
UnicastFlood=foo

