# tests/auto/InvalidValue/2.31/bad/Bridge_HairPin_baz.target
[Unit]
Description=Foo
[Bridge]
HairPin=baz

