# tests/auto/InvalidValue/2.31/bad/Bridge_Learning_bar.target
[Unit]
Description=Foo
[Bridge]
Learning=bar

