# tests/auto/InvalidValue/2.31/bad/Bridge_AllowPortToBeRoot_baz.target
[Unit]
Description=Foo
[Bridge]
AllowPortToBeRoot=baz

