# tests/auto/InvalidValue/2.31/bad/Bridge_MulticastFlood_foo.target
[Unit]
Description=Foo
[Bridge]
MulticastFlood=foo

