# tests/auto/InvalidValue/2.31/bad/Bridge_NeighborSuppression_foo.target
[Unit]
Description=Foo
[Bridge]
NeighborSuppression=foo

