# tests/auto/InvalidValue/2.31/bad/Bridge_UseBPDU_foo.target
[Unit]
Description=Foo
[Bridge]
UseBPDU=foo

