# tests/auto/InvalidValue/2.31/bad/Bridge_AllowPortToBeRoot_foo.target
[Unit]
Description=Foo
[Bridge]
AllowPortToBeRoot=foo

