# tests/auto/InvalidValue/2.31/bad/CAN_TripleSampling_foo.target
[Unit]
Description=Foo
[CAN]
TripleSampling=foo

