# tests/auto/InvalidValue/2.31/bad/CAN_RestartSec_0.5d.target
[Unit]
Description=Foo
[CAN]
RestartSec=0.5d

