# tests/auto/InvalidValue/2.31/bad/CAN_Bitrate_-1M.target
[Unit]
Description=Foo
[CAN]
Bitrate=-1M

