# tests/auto/InvalidValue/2.31/bad/DHCP_UseTimezone_bar.target
[Unit]
Description=Foo
[DHCP]
UseTimezone=bar

