# tests/auto/InvalidValue/2.31/bad/DHCP_UseMTU_bar.target
[Unit]
Description=Foo
[DHCP]
UseMTU=bar

