# tests/auto/InvalidValue/2.31/bad/DHCPServer_PoolOffset_-1.target
[Unit]
Description=Foo
[DHCPServer]
PoolOffset=-1

