# tests/auto/InvalidValue/2.31/bad/DHCPServer_EmitDNS_bar.target
[Unit]
Description=Foo
[DHCPServer]
EmitDNS=bar

