# tests/auto/InvalidValue/2.31/bad/DHCP_UseNTP_foo.target
[Unit]
Description=Foo
[DHCP]
UseNTP=foo

