# tests/auto/InvalidValue/2.31/bad/DHCPServer_PoolOffset_10000000000000.target
[Unit]
Description=Foo
[DHCPServer]
PoolOffset=10000000000000

