# tests/auto/InvalidValue/2.31/bad/DHCPServer_EmitNTP_foo.target
[Unit]
Description=Foo
[DHCPServer]
EmitNTP=foo

