# tests/auto/InvalidValue/2.31/bad/DHCP_RouteTable_0.target
[Unit]
Description=Foo
[DHCP]
RouteTable=0

