# tests/auto/InvalidValue/2.31/bad/DHCP_CriticalConnection_bar.target
[Unit]
Description=Foo
[DHCP]
CriticalConnection=bar

