# tests/auto/InvalidValue/2.31/bad/DHCPServer_EmitTimezone_baz.target
[Unit]
Description=Foo
[DHCPServer]
EmitTimezone=baz

