# tests/auto/InvalidValue/2.31/bad/DHCPServer_EmitNTP_bar.target
[Unit]
Description=Foo
[DHCPServer]
EmitNTP=bar

