# tests/auto/InvalidValue/2.31/bad/DHCPServer_EmitRouter_baz.target
[Unit]
Description=Foo
[DHCPServer]
EmitRouter=baz

