# tests/auto/InvalidValue/2.31/bad/DHCP_UseTimezone_foo.target
[Unit]
Description=Foo
[DHCP]
UseTimezone=foo

