# tests/auto/InvalidValue/2.31/bad/DHCP_ListenPort_0.target
[Unit]
Description=Foo
[DHCP]
ListenPort=0

