# tests/auto/InvalidValue/2.31/bad/DHCP_IAID_10000000000000.target
[Unit]
Description=Foo
[DHCP]
IAID=10000000000000

