# tests/auto/InvalidValue/2.31/bad/DHCPServer_PoolSize_10000000000000.target
[Unit]
Description=Foo
[DHCPServer]
PoolSize=10000000000000

