# tests/auto/InvalidValue/2.31/bad/DHCPServer_EmitDNS_foo.target
[Unit]
Description=Foo
[DHCPServer]
EmitDNS=foo

