# tests/auto/InvalidValue/2.31/bad/DHCPServer_EmitTimezone_foo.target
[Unit]
Description=Foo
[DHCPServer]
EmitTimezone=foo

