# tests/auto/InvalidValue/2.31/bad/DHCP_Anonymize_baz.target
[Unit]
Description=Foo
[DHCP]
Anonymize=baz

