# tests/auto/InvalidValue/2.31/bad/DHCP_UseTimezone_baz.target
[Unit]
Description=Foo
[DHCP]
UseTimezone=baz

