# tests/auto/InvalidValue/2.31/bad/DHCP_CriticalConnection_baz.target
[Unit]
Description=Foo
[DHCP]
CriticalConnection=baz

