# tests/auto/InvalidValue/2.31/bad/DHCP_UseRoutes_foo.target
[Unit]
Description=Foo
[DHCP]
UseRoutes=foo

