# tests/auto/InvalidValue/2.31/bad/DHCP_RapidCommit_bar.target
[Unit]
Description=Foo
[DHCP]
RapidCommit=bar

