# tests/auto/InvalidValue/2.31/bad/DHCPServer_MaxLeaseTimeSec_2potatoes.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=2potatoes

