# tests/auto/InvalidValue/2.31/bad/DHCPServer_MaxLeaseTimeSec_0.5d.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=0.5d

