# tests/auto/InvalidValue/2.31/bad/FooOverUDP_Protocol_0.target
[Unit]
Description=Foo
[FooOverUDP]
Protocol=0

