# tests/auto/InvalidValue/2.31/bad/GENEVE_TOS_256.target
[Unit]
Description=Foo
[GENEVE]
TOS=256

