# tests/auto/InvalidValue/2.31/bad/GENEVE_DestinationPort_65535.target
[Unit]
Description=Foo
[GENEVE]
DestinationPort=65535

