# tests/auto/InvalidValue/2.31/bad/GENEVE_Id_16777216.target
[Unit]
Description=Foo
[GENEVE]
Id=16777216

