# tests/auto/InvalidValue/2.31/bad/GENEVE_Remote_abc.target
[Unit]
Description=Foo
[GENEVE]
Remote=abc

