# tests/auto/InvalidValue/2.31/bad/GENEVE_Remote_255.255.255.256.target
[Unit]
Description=Foo
[GENEVE]
Remote=255.255.255.256

