# tests/auto/InvalidValue/2.31/bad/GENEVE_UDPChecksum_baz.target
[Unit]
Description=Foo
[GENEVE]
UDPChecksum=baz

