# tests/auto/InvalidValue/2.31/bad/GENEVE_Remote_0.0.0..target
[Unit]
Description=Foo
[GENEVE]
Remote=0.0.0.

