# tests/auto/InvalidValue/2.31/bad/IPv6PrefixDelegation_DNSLifetimeSec_3decades.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
DNSLifetimeSec=3decades

