# tests/auto/InvalidValue/2.31/bad/IPv6AcceptRA_UseDomains_foo.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseDomains=foo

