# tests/auto/InvalidValue/2.31/bad/IPv6Prefix_AddressAutoconfiguration_baz.target
[Unit]
Description=Foo
[IPv6Prefix]
AddressAutoconfiguration=baz

