# tests/auto/InvalidValue/2.31/bad/IPv6AcceptRA_UseAutonomousPrefix_bar.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseAutonomousPrefix=bar

