# tests/auto/InvalidValue/2.31/bad/IPv6Prefix_ValidLifetimeSec_True.target
[Unit]
Description=Foo
[IPv6Prefix]
ValidLifetimeSec=True

