# tests/auto/InvalidValue/2.31/bad/IPv6PrefixDelegation_DNSLifetimeSec_True.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
DNSLifetimeSec=True

