# tests/auto/InvalidValue/2.31/bad/IPv6PrefixDelegation_EmitDNS_foo.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
EmitDNS=foo

