# tests/auto/InvalidValue/2.31/bad/IPv6AcceptRA_UseDNS_foo.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseDNS=foo

