# tests/auto/InvalidValue/2.31/bad/Journal_LineMax_78T.target
[Unit]
Description=Foo
[Journal]
LineMax=78T

