# tests/auto/InvalidValue/2.31/bad/Journal_ForwardToConsole_bar.target
[Unit]
Description=Foo
[Journal]
ForwardToConsole=bar

