# tests/auto/InvalidValue/2.31/bad/Journal_RuntimeKeepFree_10000000000000T.target
[Unit]
Description=Foo
[Journal]
RuntimeKeepFree=10000000000000T

