# tests/auto/InvalidValue/2.31/bad/Journal_MaxFileSec_-1day.target
[Unit]
Description=Foo
[Journal]
MaxFileSec=-1day

