# tests/auto/InvalidValue/2.31/bad/Journal_SyncIntervalSec_2potatoes.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=2potatoes

