# tests/auto/InvalidValue/2.31/bad/Journal_RateLimitBurst_10000000000000.target
[Unit]
Description=Foo
[Journal]
RateLimitBurst=10000000000000

