# tests/auto/InvalidValue/2.31/bad/Journal_Seal_bar.target
[Unit]
Description=Foo
[Journal]
Seal=bar

