# tests/auto/InvalidValue/2.31/bad/Journal_ForwardToSyslog_baz.target
[Unit]
Description=Foo
[Journal]
ForwardToSyslog=baz

