# tests/auto/InvalidValue/2.31/bad/Journal_RateLimitBurst_-1.target
[Unit]
Description=Foo
[Journal]
RateLimitBurst=-1

