# tests/auto/InvalidValue/2.31/bad/Journal_RateLimitIntervalSec_-1h.target
[Unit]
Description=Foo
[Journal]
RateLimitIntervalSec=-1h

