# tests/auto/InvalidValue/2.31/bad/Journal_Seal_baz.target
[Unit]
Description=Foo
[Journal]
Seal=baz

