# tests/auto/InvalidValue/2.31/bad/Journal_SystemMaxFileSize_-1K.target
[Unit]
Description=Foo
[Journal]
SystemMaxFileSize=-1K

