# tests/auto/InvalidValue/2.31/bad/Journal_SystemKeepFree_-1P.target
[Unit]
Description=Foo
[Journal]
SystemKeepFree=-1P

