# tests/auto/InvalidValue/2.31/bad/Journal_RuntimeMaxFiles_10000000000000K.target
[Unit]
Description=Foo
[Journal]
RuntimeMaxFiles=10000000000000K

