# tests/auto/InvalidValue/2.31/bad/Journal_ForwardToKMsg_bar.target
[Unit]
Description=Foo
[Journal]
ForwardToKMsg=bar

