# tests/auto/InvalidValue/2.31/bad/Journal_ReadKMsg_baz.target
[Unit]
Description=Foo
[Journal]
ReadKMsg=baz

