# tests/auto/InvalidValue/2.31/bad/Journal_MaxFileSec_-1month.target
[Unit]
Description=Foo
[Journal]
MaxFileSec=-1month

