# tests/auto/InvalidValue/2.31/bad/Journal_SystemMaxFiles_10000000000000P.target
[Unit]
Description=Foo
[Journal]
SystemMaxFiles=10000000000000P

