# tests/auto/InvalidValue/2.31/bad/Journal_ForwardToSyslog_bar.target
[Unit]
Description=Foo
[Journal]
ForwardToSyslog=bar

