# tests/auto/InvalidValue/2.31/bad/Journal_ForwardToKMsg_baz.target
[Unit]
Description=Foo
[Journal]
ForwardToKMsg=baz

