# tests/auto/InvalidValue/2.31/bad/Journal_SystemMaxUse_10000000000000K.target
[Unit]
Description=Foo
[Journal]
SystemMaxUse=10000000000000K

