# tests/auto/InvalidValue/2.31/bad/Journal_Seal_foo.target
[Unit]
Description=Foo
[Journal]
Seal=foo

