# tests/auto/InvalidValue/2.31/bad/Journal_MaxRetentionSec_10000000000000h.target
[Unit]
Description=Foo
[Journal]
MaxRetentionSec=10000000000000h

