# tests/auto/InvalidValue/2.31/bad/Journal_SystemMaxFiles_-1T.target
[Unit]
Description=Foo
[Journal]
SystemMaxFiles=-1T

