# tests/auto/InvalidValue/2.31/bad/Journal_RateLimitIntervalSec_10000000000000ms.target
[Unit]
Description=Foo
[Journal]
RateLimitIntervalSec=10000000000000ms

