# tests/auto/InvalidValue/2.31/bad/Journal_ForwardToWall_foo.target
[Unit]
Description=Foo
[Journal]
ForwardToWall=foo

