# tests/auto/InvalidValue/2.31/bad/Journal_ForwardToWall_baz.target
[Unit]
Description=Foo
[Journal]
ForwardToWall=baz

