# tests/auto/InvalidValue/2.31/bad/Journal_RuntimeMaxFiles_-1E.target
[Unit]
Description=Foo
[Journal]
RuntimeMaxFiles=-1E

