# tests/auto/InvalidValue/2.31/bad/Journal_RuntimeMaxFileSize_10000000000000E.target
[Unit]
Description=Foo
[Journal]
RuntimeMaxFileSize=10000000000000E

