# tests/auto/InvalidValue/2.31/bad/L2TP_Remote_1.1.1.1---1.target
[Unit]
Description=Foo
[L2TP]
Remote=1.1.1.1/-1

