# tests/auto/InvalidValue/2.31/bad/L2TP_UDPChecksum_baz.target
[Unit]
Description=Foo
[L2TP]
UDPChecksum=baz

