# tests/auto/InvalidValue/2.31/bad/L2TP_Remote_0.0.0..target
[Unit]
Description=Foo
[L2TP]
Remote=0.0.0.

