# tests/auto/InvalidValue/2.31/bad/L2TP_Remote_255.256.255.255.target
[Unit]
Description=Foo
[L2TP]
Remote=255.256.255.255

