# tests/auto/InvalidValue/2.31/bad/L2TP_TunnelId_0.target
[Unit]
Description=Foo
[L2TP]
TunnelId=0

