# tests/auto/InvalidValue/2.31/bad/L2TPSession_SessionId_4294967296.target
[Unit]
Description=Foo
[L2TPSession]
Name=abc
SessionId=4294967296

