# tests/auto/InvalidValue/2.31/bad/L2TP_DestinationPort_0.target
[Unit]
Description=Foo
[L2TP]
DestinationPort=0

