# tests/auto/InvalidValue/2.31/bad/Link_AllMulticast_foo.target
[Unit]
Description=Foo
[Link]
AllMulticast=foo

