# tests/auto/InvalidValue/2.31/bad/Link_MTUBytes_-1M.target
[Unit]
Description=Foo
[Link]
MTUBytes=-1M

