# tests/auto/InvalidValue/2.31/bad/Link_Unmanaged_bar.target
[Unit]
Description=Foo
[Link]
Unmanaged=bar

