# tests/auto/InvalidValue/2.31/bad/Link_Multicast_baz.target
[Unit]
Description=Foo
[Link]
Multicast=baz

