# tests/auto/InvalidValue/2.31/bad/Link_RequiredForOnline_foo.target
[Unit]
Description=Foo
[Link]
RequiredForOnline=foo

