# tests/auto/InvalidValue/2.31/bad/Link_Multicast_foo.target
[Unit]
Description=Foo
[Link]
Multicast=foo

