# tests/auto/InvalidValue/2.31/bad/Link_AllMulticast_baz.target
[Unit]
Description=Foo
[Link]
AllMulticast=baz

