# tests/auto/InvalidValue/2.31/bad/Login_InhibitorsMax_-1T.target
[Unit]
Description=Foo
[Login]
InhibitorsMax=-1T

