# tests/auto/InvalidValue/2.31/bad/Login_HoldoffTimeoutSec_True.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=True

