# tests/auto/InvalidValue/2.31/bad/Login_SuspendKeyIgnoreInhibited_baz.target
[Unit]
Description=Foo
[Login]
SuspendKeyIgnoreInhibited=baz

