# tests/auto/InvalidValue/2.31/bad/Login_PowerKeyIgnoreInhibited_baz.target
[Unit]
Description=Foo
[Login]
PowerKeyIgnoreInhibited=baz

