# tests/auto/InvalidValue/2.31/bad/Login_IdleActionSec_True.target
[Unit]
Description=Foo
[Login]
IdleActionSec=True

