# tests/auto/InvalidValue/2.31/bad/Login_LidSwitchIgnoreInhibited_baz.target
[Unit]
Description=Foo
[Login]
LidSwitchIgnoreInhibited=baz

