# tests/auto/InvalidValue/2.31/bad/Login_IdleActionSec_1p.target
[Unit]
Description=Foo
[Login]
IdleActionSec=1p

