# tests/auto/InvalidValue/2.31/bad/Login_ReserveVT_10000000000000.target
[Unit]
Description=Foo
[Login]
ReserveVT=10000000000000

