# tests/auto/InvalidValue/2.31/bad/Login_UserStopDelaySec_3decades.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=3decades

