# tests/auto/InvalidValue/2.31/bad/Login_PowerKeyIgnoreInhibited_bar.target
[Unit]
Description=Foo
[Login]
PowerKeyIgnoreInhibited=bar

