# tests/auto/InvalidValue/2.31/bad/Login_RuntimeDirectorySize_10000000000000K.target
[Unit]
Description=Foo
[Login]
RuntimeDirectorySize=10000000000000K

