# tests/auto/InvalidValue/2.31/bad/Login_KillUserProcesses_foo.target
[Unit]
Description=Foo
[Login]
KillUserProcesses=foo

