# tests/auto/InvalidValue/2.31/bad/Login_RemoveIPC_baz.target
[Unit]
Description=Foo
[Login]
RemoveIPC=baz

