# tests/auto/InvalidValue/2.31/bad/Login_HibernateKeyIgnoreInhibited_foo.target
[Unit]
Description=Foo
[Login]
HibernateKeyIgnoreInhibited=foo

