# tests/auto/InvalidValue/2.31/bad/Login_ReserveVT_-1.target
[Unit]
Description=Foo
[Login]
ReserveVT=-1

