# tests/auto/InvalidValue/2.31/bad/Login_HoldoffTimeoutSec_3decades.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=3decades

