# tests/auto/InvalidValue/2.31/bad/Login_UserStopDelaySec_True.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=True

