# tests/auto/InvalidValue/2.31/bad/Login_LidSwitchIgnoreInhibited_foo.target
[Unit]
Description=Foo
[Login]
LidSwitchIgnoreInhibited=foo

