# tests/auto/InvalidValue/2.31/bad/Login_RemoveIPC_bar.target
[Unit]
Description=Foo
[Login]
RemoveIPC=bar

