# tests/auto/InvalidValue/2.31/bad/Login_PowerKeyIgnoreInhibited_foo.target
[Unit]
Description=Foo
[Login]
PowerKeyIgnoreInhibited=foo

