# tests/auto/InvalidValue/2.31/bad/Login_HoldoffTimeoutSec_0.5d.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=0.5d

