# tests/auto/InvalidValue/2.31/bad/Login_HoldoffTimeoutSec_2potatoes.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=2potatoes

