# tests/auto/InvalidValue/2.31/bad/Login_HibernateKeyIgnoreInhibited_baz.target
[Unit]
Description=Foo
[Login]
HibernateKeyIgnoreInhibited=baz

