# tests/auto/InvalidValue/2.31/bad/Manager_DefaultTasksAccounting_bar.target
[Unit]
Description=Foo
[Manager]
DefaultTasksAccounting=bar

