# tests/auto/InvalidValue/2.31/bad/Manager_RuntimeWatchdogSec_True.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=True

