# tests/auto/InvalidValue/2.31/bad/Manager_DefaultRestartSec_3decades.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=3decades

