# tests/auto/InvalidValue/2.31/bad/Manager_DumpCore_bar.target
[Unit]
Description=Foo
[Manager]
DumpCore=bar

