# tests/auto/InvalidValue/2.31/bad/Manager_DefaultTasksMax_10000000000000%25.target
[Unit]
Description=Foo
[Manager]
DefaultTasksMax=10000000000000%

