# tests/auto/InvalidValue/2.31/bad/Manager_ShutdownWatchdogSec_2potatoes.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=2potatoes

