# tests/auto/InvalidValue/2.31/bad/Manager_DefaultRestartSec_0.5d.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=0.5d

