# tests/auto/InvalidValue/2.31/bad/Manager_DefaultRestartSec_True.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=True

