# tests/auto/InvalidValue/2.31/bad/Manager_TimerSlackNSec_2potatoes.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=2potatoes

