# tests/auto/InvalidValue/2.31/bad/Manager_WatchdogDevice_--tmp.target
[Unit]
Description=Foo
[Manager]
WatchdogDevice=/tmp

