# tests/auto/InvalidValue/2.31/bad/Manager_DefaultTimerAccuracySec_2potatoes.target
[Unit]
Description=Foo
[Manager]
DefaultTimerAccuracySec=2potatoes

