# tests/auto/InvalidValue/2.31/bad/Manager_DefaultTimerAccuracySec_True.target
[Unit]
Description=Foo
[Manager]
DefaultTimerAccuracySec=True

