# tests/auto/InvalidValue/2.31/bad/Manager_SystemCallArchitectures_1p.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=1p

