# tests/auto/InvalidValue/2.31/bad/Manager_DefaultBlockIOAccounting_bar.target
[Unit]
Description=Foo
[Manager]
DefaultBlockIOAccounting=bar

