# tests/auto/InvalidValue/2.31/bad/Manager_DefaultCPUAccounting_bar.target
[Unit]
Description=Foo
[Manager]
DefaultCPUAccounting=bar

