# tests/auto/InvalidValue/2.31/bad/Manager_DefaultTimeoutStopSec_0.5d.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=0.5d

