# tests/auto/InvalidValue/2.31/bad/Manager_RuntimeWatchdogSec_3decades.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=3decades

