# tests/auto/InvalidValue/2.31/bad/Manager_DefaultIOAccounting_bar.target
[Unit]
Description=Foo
[Manager]
DefaultIOAccounting=bar

