# tests/auto/InvalidValue/2.31/bad/Manager_DefaultMemoryAccounting_foo.target
[Unit]
Description=Foo
[Manager]
DefaultMemoryAccounting=foo

