# tests/auto/InvalidValue/2.31/bad/Manager_ShowStatus_bar.target
[Unit]
Description=Foo
[Manager]
ShowStatus=bar

