# tests/auto/InvalidValue/2.31/bad/Manager_DefaultCPUAccounting_baz.target
[Unit]
Description=Foo
[Manager]
DefaultCPUAccounting=baz

