# tests/auto/InvalidValue/2.31/bad/Manager_DefaultBlockIOAccounting_foo.target
[Unit]
Description=Foo
[Manager]
DefaultBlockIOAccounting=foo

