# tests/auto/InvalidValue/2.31/bad/Manager_DefaultRestartSec_2potatoes.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=2potatoes

