# tests/auto/InvalidValue/2.31/bad/Manager_CrashReboot_foo.target
[Unit]
Description=Foo
[Manager]
CrashReboot=foo

