# tests/auto/InvalidValue/2.31/bad/Manager_DefaultIPAccounting_bar.target
[Unit]
Description=Foo
[Manager]
DefaultIPAccounting=bar

