# tests/auto/InvalidValue/2.31/bad/Manager_TimerSlackNSec_True.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=True

