# tests/auto/InvalidValue/2.31/bad/Manager_ShowStatus_baz.target
[Unit]
Description=Foo
[Manager]
ShowStatus=baz

