# tests/auto/InvalidValue/2.31/bad/Manager_CPUAffinity_10000000000000.target
[Unit]
Description=Foo
[Manager]
CPUAffinity=10000000000000

