# tests/auto/InvalidValue/2.31/bad/Manager_LogColor_foo.target
[Unit]
Description=Foo
[Manager]
LogColor=foo

