# tests/auto/InvalidValue/2.31/bad/Manager_DefaultEnvironment_1.target
[Unit]
Description=Foo
[Manager]
DefaultEnvironment=1

