# tests/auto/InvalidValue/2.31/bad/Manager_LogColor_bar.target
[Unit]
Description=Foo
[Manager]
LogColor=bar

