# tests/auto/InvalidValue/2.31/bad/Manager_ShutdownWatchdogSec_0.5d.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=0.5d

