# tests/auto/InvalidValue/2.31/bad/Manager_SystemCallArchitectures_0.5d.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=0.5d

