# tests/auto/InvalidValue/2.31/bad/Manager_DefaultStartLimitIntervalSec_True.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=True

