# tests/auto/InvalidValue/2.31/bad/Match_MACAddress_01%3A23%3A45%3A67%3A89.target
[Unit]
Description=Foo
[Match]
MACAddress=01:23:45:67:89

