# tests/auto/InvalidValue/2.31/bad/Mount_MemoryMax_10000000000000K.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
MemoryMax=10000000000000K

