# tests/auto/InvalidValue/2.31/bad/Mount_LimitRTPRIO_10000000000000.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
LimitRTPRIO=10000000000000

