# tests/auto/InvalidValue/2.31/bad/Mount_BlockIOWriteBandwidthMax_--tmp--dev%2010000000000000T.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
BlockIOWriteBandwidthMax=/tmp/dev 10000000000000T

