# tests/auto/InvalidValue/2.31/bad/Mount_LimitNICE_-1.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
LimitNICE=-1

