# tests/auto/InvalidValue/2.31/bad/Mount_BlockIODeviceWeight_--tmp%2010001.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
BlockIODeviceWeight=/tmp 10001

