# tests/auto/InvalidValue/2.31/bad/Mount_LimitDATA_10000000000000.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
LimitDATA=10000000000000

