# tests/auto/InvalidValue/2.31/bad/Mount_BlockIOReadBandwidthMax_--devvvv%2010000000000000G.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
BlockIOReadBandwidthMax=/devvvv 10000000000000G

