# tests/auto/InvalidValue/2.31/bad/Mount_BlockIOWriteBandwidthMax_--tmp%20-1M.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
BlockIOWriteBandwidthMax=/tmp -1M

