# tests/auto/InvalidValue/2.31/bad/Mount_BlockIOWriteBandwidthMax_--tmp--dev%20-1T.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
BlockIOWriteBandwidthMax=/tmp/dev -1T

