# tests/auto/InvalidValue/2.31/bad/Mount_BlockIOReadBandwidthMax_dev--foo%20-1K.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
BlockIOReadBandwidthMax=dev/foo -1K

