# tests/auto/InvalidValue/2.31/bad/Mount_Where_this.isnt--gonna-work_here.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
Where=this.isnt--gonna-work_here

