# tests/auto/InvalidValue/2.31/bad/Mount_LimitNOFILE_10000000000000.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
LimitNOFILE=10000000000000

