# tests/auto/InvalidValue/2.31/bad/Mount_BlockIOReadBandwidthMax_--tmp--dev%20-1G.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
BlockIOReadBandwidthMax=/tmp/dev -1G

