# tests/auto/InvalidValue/2.31/bad/Mount_LimitFSIZE_-1.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
LimitFSIZE=-1

