# tests/auto/InvalidValue/2.31/bad/Neighbor_Address_0--abc.target
[Unit]
Description=Foo
[Neighbor]
Address=0/abc

