# tests/auto/InvalidValue/2.31/bad/Neighbor_Address_abc.target
[Unit]
Description=Foo
[Neighbor]
Address=abc

