# tests/auto/InvalidValue/2.31/bad/NetDev_MTUBytes_-1K.target
[Unit]
Description=Foo
[NetDev]
Name=abc
MTUBytes=-1K

