# tests/auto/InvalidValue/2.31/bad/NetDev_MTUBytes_10000000000000K.target
[Unit]
Description=Foo
[NetDev]
Name=abc
MTUBytes=10000000000000K

