# tests/auto/InvalidValue/2.31/bad/Network_IPv4ProxyARP_bar.target
[Unit]
Description=Foo
[Network]
IPv4ProxyARP=bar

