# tests/auto/InvalidValue/2.31/bad/Network_IPv6ProxyNDPAddress_abc.target
[Unit]
Description=Foo
[Network]
IPv6ProxyNDPAddress=abc

