# tests/auto/InvalidValue/2.31/bad/Network_DNSDefaultRoute_bar.target
[Unit]
Description=Foo
[Network]
DNSDefaultRoute=bar

