# tests/auto/InvalidValue/2.31/bad/Network_IPv4LLRoute_foo.target
[Unit]
Description=Foo
[Network]
IPv4LLRoute=foo

