# tests/auto/InvalidValue/2.31/bad/Network_IPv6Token_0--abc.target
[Unit]
Description=Foo
[Network]
IPv6Token=0/abc

