# tests/auto/InvalidValue/2.31/bad/Network_DHCPServer_foo.target
[Unit]
Description=Foo
[Network]
DHCPServer=foo

