# tests/auto/InvalidValue/2.31/bad/Network_IPv6ProxyNDPAddress_0--abc.target
[Unit]
Description=Foo
[Network]
IPv6ProxyNDPAddress=0/abc

