# tests/auto/InvalidValue/2.31/bad/Network_DHCPServer_baz.target
[Unit]
Description=Foo
[Network]
DHCPServer=baz

