# tests/auto/InvalidValue/2.31/bad/Network_NTP_0.0.0..target
[Unit]
Description=Foo
[Network]
NTP=0.0.0.

