# tests/auto/InvalidValue/2.31/bad/Network_ConfigureWithoutCarrier_baz.target
[Unit]
Description=Foo
[Network]
ConfigureWithoutCarrier=baz

