# tests/auto/InvalidValue/2.31/bad/Network_DNS_abc.target
[Unit]
Description=Foo
[Network]
DNS=abc

