# tests/auto/InvalidValue/2.31/bad/Network_IPv6HopLimit_-1.target
[Unit]
Description=Foo
[Network]
IPv6HopLimit=-1

