# tests/auto/InvalidValue/2.31/bad/Network_Address_255.255.255.256.target
[Unit]
Description=Foo
[Network]
Address=255.255.255.256

