# tests/auto/InvalidValue/2.31/bad/Network_IPv6DuplicateAddressDetection_-1.target
[Unit]
Description=Foo
[Network]
IPv6DuplicateAddressDetection=-1

