# tests/auto/InvalidValue/2.31/bad/Network_DNSDefaultRoute_baz.target
[Unit]
Description=Foo
[Network]
DNSDefaultRoute=baz

