# tests/auto/InvalidValue/2.31/bad/Network_Address_abc.target
[Unit]
Description=Foo
[Network]
Address=abc

