# tests/auto/InvalidValue/2.31/bad/Path_DirectoryNotEmpty_..--..--.target
[Unit]
Description=Foo
[Path]
DirectoryNotEmpty=../../

