# tests/auto/InvalidValue/2.31/bad/Path_PathExists_True.target
[Unit]
Description=Foo
[Path]
PathExists=True

