# tests/auto/InvalidValue/2.31/bad/Path_PathExists_1.target
[Unit]
Description=Foo
[Path]
PathExists=1

