# tests/auto/InvalidValue/2.31/bad/Path_PathModified_1.target
[Unit]
Description=Foo
[Path]
PathModified=1

