# tests/auto/InvalidValue/2.31/bad/Path_PathModified_..--..--.target
[Unit]
Description=Foo
[Path]
PathModified=../../

