# tests/auto/InvalidValue/2.31/bad/Path_DirectoryMode_True.target
[Unit]
Description=Foo
[Path]
DirectoryMode=True

