# tests/auto/InvalidValue/2.31/bad/Path_DirectoryNotEmpty_tmp.target
[Unit]
Description=Foo
[Path]
DirectoryNotEmpty=tmp

