# tests/auto/InvalidValue/2.31/bad/Path_MakeDirectory_foo.target
[Unit]
Description=Foo
[Path]
MakeDirectory=foo

