# tests/auto/InvalidValue/2.31/bad/Path_DirectoryMode_888.target
[Unit]
Description=Foo
[Path]
DirectoryMode=888

