# tests/auto/InvalidValue/2.31/bad/Path_PathChanged_True.target
[Unit]
Description=Foo
[Path]
PathChanged=True

