# tests/auto/InvalidValue/2.31/bad/Path_DirectoryNotEmpty_True.target
[Unit]
Description=Foo
[Path]
DirectoryNotEmpty=True

