# tests/auto/InvalidValue/2.31/bad/Route_QuickAck_baz.target
[Unit]
Description=Foo
[Route]
QuickAck=baz

