# tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_To_0.0.0..target
[Unit]
Description=Foo
[RoutingPolicyRule]
To=0.0.0.

