# tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_From_256.255.255.255.target
[Unit]
Description=Foo
[RoutingPolicyRule]
From=256.255.255.255

