# tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_SourcePort_-1%2010000000000000.target
[Unit]
Description=Foo
[RoutingPolicyRule]
SourcePort=-1 10000000000000

