# tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_InvertRule_baz.target
[Unit]
Description=Foo
[RoutingPolicyRule]
InvertRule=baz

