# tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_Table_4294967296.target
[Unit]
Description=Foo
[RoutingPolicyRule]
Table=4294967296

