# tests/auto/InvalidValue/2.31/bad/Route_GatewayOnLink_bar.target
[Unit]
Description=Foo
[Route]
GatewayOnLink=bar

