# tests/auto/InvalidValue/2.31/bad/Route_InitialCongestionWindow_4294967296.target
[Unit]
Description=Foo
[Route]
InitialCongestionWindow=4294967296

