# tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_DestinationPort_-1%20-1.target
[Unit]
Description=Foo
[RoutingPolicyRule]
DestinationPort=-1 -1

