# tests/auto/InvalidValue/2.31/bad/Route_MTUBytes_10000000000000K.target
[Unit]
Description=Foo
[Route]
MTUBytes=10000000000000K

