# tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_Priority_10000000000000.target
[Unit]
Description=Foo
[RoutingPolicyRule]
Priority=10000000000000

