# tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_From_abc.target
[Unit]
Description=Foo
[RoutingPolicyRule]
From=abc

