# tests/auto/InvalidValue/2.31/bad/Route_Metric_10000000000000.target
[Unit]
Description=Foo
[Route]
Metric=10000000000000

