# tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_To_255.256.255.255.target
[Unit]
Description=Foo
[RoutingPolicyRule]
To=255.256.255.255

