# tests/auto/InvalidValue/2.31/bad/Route_GatewayOnLink_baz.target
[Unit]
Description=Foo
[Route]
GatewayOnLink=baz

