# tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_FirewallMark_0.target
[Unit]
Description=Foo
[RoutingPolicyRule]
FirewallMark=0

