# tests/auto/InvalidValue/2.31/bad/Route_InitialCongestionWindow_0.target
[Unit]
Description=Foo
[Route]
InitialCongestionWindow=0

