# tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_Priority_-1.target
[Unit]
Description=Foo
[RoutingPolicyRule]
Priority=-1

