# tests/auto/InvalidValue/2.31/bad/RoutingPolicyRule_To_abc.target
[Unit]
Description=Foo
[RoutingPolicyRule]
To=abc

