# tests/auto/InvalidValue/2.31/bad/Service_User_True.target
[Unit]
Description=Foo
[Service]
User=True

