# tests/auto/InvalidValue/2.31/bad/Service_IOSchedulingPriority_8.target
[Unit]
Description=Foo
[Service]
IOSchedulingPriority=8

