# tests/auto/InvalidValue/2.31/bad/Service_MemoryMax_-1G.target
[Unit]
Description=Foo
[Service]
MemoryMax=-1G

