# tests/auto/InvalidValue/2.31/bad/Service_EnvironmentFile_True.target
[Unit]
Description=Foo
[Service]
EnvironmentFile=True

