# tests/auto/InvalidValue/2.31/bad/Service_LimitCORE_-1.target
[Unit]
Description=Foo
[Service]
LimitCORE=-1

