# tests/auto/InvalidValue/2.31/bad/Service_StartupCPUShares_262145.target
[Unit]
Description=Foo
[Service]
StartupCPUShares=262145

