# tests/auto/InvalidValue/2.31/bad/Service_PrivateTmp_baz.target
[Unit]
Description=Foo
[Service]
PrivateTmp=baz

