# tests/auto/InvalidValue/2.31/bad/Service_UMask_-1.target
[Unit]
Description=Foo
[Service]
UMask=-1

