# tests/auto/InvalidValue/2.31/bad/Service_BlockIOWriteBandwidthMax_dev--foo%2010000000000000G.target
[Unit]
Description=Foo
[Service]
BlockIOWriteBandwidthMax=dev/foo 10000000000000G

