# tests/auto/InvalidValue/2.31/bad/Service_TimeoutStartSec_0.5d.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=0.5d

