# tests/auto/InvalidValue/2.31/bad/Service_TasksMax_0.target
[Unit]
Description=Foo
[Service]
TasksMax=0

