# tests/auto/InvalidValue/2.31/bad/Service_TasksMax_10000000000000.target
[Unit]
Description=Foo
[Service]
TasksMax=10000000000000

