# tests/auto/InvalidValue/2.31/bad/Service_PrivateNetwork_baz.target
[Unit]
Description=Foo
[Service]
PrivateNetwork=baz

