# tests/auto/InvalidValue/2.31/bad/Service_TTYVTDisallocate_bar.target
[Unit]
Description=Foo
[Service]
TTYVTDisallocate=bar

