# tests/auto/InvalidValue/2.31/bad/Service_IOAccounting_bar.target
[Unit]
Description=Foo
[Service]
IOAccounting=bar

