# tests/auto/InvalidValue/2.31/bad/Service_LimitCPU_-1.target
[Unit]
Description=Foo
[Service]
LimitCPU=-1

