# tests/auto/InvalidValue/2.31/bad/Service_MemorySwapMax_-1T.target
[Unit]
Description=Foo
[Service]
MemorySwapMax=-1T

