# tests/auto/InvalidValue/2.31/bad/Service_SendSIGHUP_bar.target
[Unit]
Description=Foo
[Service]
SendSIGHUP=bar

