# tests/auto/InvalidValue/2.31/bad/Service_RestrictSUIDSGID_baz.target
[Unit]
Description=Foo
[Service]
RestrictSUIDSGID=baz

