# tests/auto/InvalidValue/2.31/bad/Service_SupplementaryGroups_3.target
[Unit]
Description=Foo
[Service]
SupplementaryGroups=3

