# tests/auto/InvalidValue/2.31/bad/Service_RuntimeMaxSec_2potatoes.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=2potatoes

