# tests/auto/InvalidValue/2.31/bad/Service_StartLimitBurst_10000000000000.target
[Unit]
Description=Foo
[Service]
StartLimitBurst=10000000000000

