# tests/auto/InvalidValue/2.31/bad/Service_RestartSec_1p.target
[Unit]
Description=Foo
[Service]
RestartSec=1p

