# tests/auto/InvalidValue/2.31/bad/Service_LimitMEMLOCK_10000000000000.target
[Unit]
Description=Foo
[Service]
LimitMEMLOCK=10000000000000

