# tests/auto/InvalidValue/2.31/bad/Service_Environment_True.target
[Unit]
Description=Foo
[Service]
Environment=True

