# tests/auto/InvalidValue/2.31/bad/Service_BlockIOWriteBandwidthMax_--tmp%20-1K.target
[Unit]
Description=Foo
[Service]
BlockIOWriteBandwidthMax=/tmp -1K

