# tests/auto/InvalidValue/2.31/bad/Service_MemoryAccounting_foo.target
[Unit]
Description=Foo
[Service]
MemoryAccounting=foo

