# tests/auto/InvalidValue/2.31/bad/Service_IgnoreSIGPIPE_foo.target
[Unit]
Description=Foo
[Service]
IgnoreSIGPIPE=foo

