# tests/auto/InvalidValue/2.31/bad/Service_RuntimeMaxSec_0.5d.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=0.5d

