# tests/auto/InvalidValue/2.31/bad/Service_LimitFSIZE_-1.target
[Unit]
Description=Foo
[Service]
LimitFSIZE=-1

