# tests/auto/InvalidValue/2.31/bad/Service_TasksAccounting_baz.target
[Unit]
Description=Foo
[Service]
TasksAccounting=baz

