# tests/auto/InvalidValue/2.31/bad/Service_WatchdogSec_2potatoes.target
[Unit]
Description=Foo
[Service]
WatchdogSec=2potatoes

