# tests/auto/InvalidValue/2.31/bad/Service_MemoryAccounting_baz.target
[Unit]
Description=Foo
[Service]
MemoryAccounting=baz

