# tests/auto/InvalidValue/2.31/bad/Service_MemoryHigh_-1M.target
[Unit]
Description=Foo
[Service]
MemoryHigh=-1M

