# tests/auto/InvalidValue/2.31/bad/Service_IgnoreSIGPIPE_baz.target
[Unit]
Description=Foo
[Service]
IgnoreSIGPIPE=baz

