# tests/auto/InvalidValue/2.31/bad/Service_IOSchedulingPriority_-1.target
[Unit]
Description=Foo
[Service]
IOSchedulingPriority=-1

