# tests/auto/InvalidValue/2.31/bad/Service_Environment_a%3Eb.target
[Unit]
Description=Foo
[Service]
Environment=a>b

