# tests/auto/InvalidValue/2.31/bad/Service_TimerSlackNSec_True.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=True

