# tests/auto/InvalidValue/2.31/bad/Service_BlockStartupIOWeight_0.target
[Unit]
Description=Foo
[Service]
BlockStartupIOWeight=0

