# tests/auto/InvalidValue/2.31/bad/Service_TimeoutStopSec_2potatoes.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=2potatoes

