# tests/auto/InvalidValue/2.31/bad/Service_SendSIGHUP_foo.target
[Unit]
Description=Foo
[Service]
SendSIGHUP=foo

