# tests/auto/InvalidValue/2.31/bad/Service_MemorySwapMax_10000000000000G.target
[Unit]
Description=Foo
[Service]
MemorySwapMax=10000000000000G

