# tests/auto/InvalidValue/2.31/bad/Service_LimitFSIZE_10000000000000.target
[Unit]
Description=Foo
[Service]
LimitFSIZE=10000000000000

