# tests/auto/InvalidValue/2.31/bad/Service_LimitNICE_10000000000000.target
[Unit]
Description=Foo
[Service]
LimitNICE=10000000000000

