# tests/auto/InvalidValue/2.31/bad/Service_CPUQuota_-1%25.target
[Unit]
Description=Foo
[Service]
CPUQuota=-1%

