# tests/auto/InvalidValue/2.31/bad/Service_LimitRSS_-1.target
[Unit]
Description=Foo
[Service]
LimitRSS=-1

