# tests/auto/InvalidValue/2.31/bad/Service_StartLimitIntervalSec_0.5d.target
[Unit]
Description=Foo
[Service]
StartLimitIntervalSec=0.5d

