# tests/auto/InvalidValue/2.31/bad/Service_MemoryLimit_10000000000000B.target
[Unit]
Description=Foo
[Service]
MemoryLimit=10000000000000B

