# tests/auto/InvalidValue/2.31/bad/Service_CPUShares_262145.target
[Unit]
Description=Foo
[Service]
CPUShares=262145

