# tests/auto/InvalidValue/2.31/bad/Service_RootDirectoryStartOnly_bar.target
[Unit]
Description=Foo
[Service]
RootDirectoryStartOnly=bar

