# tests/auto/InvalidValue/2.31/bad/Service_LimitMSGQUEUE_10000000000000.target
[Unit]
Description=Foo
[Service]
LimitMSGQUEUE=10000000000000

