# tests/auto/InvalidValue/2.31/bad/Service_BlockIODeviceWeight_--devvvv%200.target
[Unit]
Description=Foo
[Service]
BlockIODeviceWeight=/devvvv 0

