# tests/auto/InvalidValue/2.31/bad/Service_RemainAfterExit_bar.target
[Unit]
Description=Foo
[Service]
RemainAfterExit=bar

