# tests/auto/InvalidValue/2.31/bad/Service_RootDirectoryStartOnly_foo.target
[Unit]
Description=Foo
[Service]
RootDirectoryStartOnly=foo

