# tests/auto/InvalidValue/2.31/bad/Service_LimitRTTIME_-1.target
[Unit]
Description=Foo
[Service]
LimitRTTIME=-1

