# tests/auto/InvalidValue/2.31/bad/Service_LimitRTTIME_10000000000000.target
[Unit]
Description=Foo
[Service]
LimitRTTIME=10000000000000

