# tests/auto/InvalidValue/2.31/bad/Service_TTYReset_foo.target
[Unit]
Description=Foo
[Service]
TTYReset=foo

