# tests/auto/InvalidValue/2.31/bad/Service_IOAccounting_foo.target
[Unit]
Description=Foo
[Service]
IOAccounting=foo

