# tests/auto/InvalidValue/2.31/bad/Service_WatchdogSec_0.5d.target
[Unit]
Description=Foo
[Service]
WatchdogSec=0.5d

