# tests/auto/InvalidValue/2.31/bad/Service_SendSIGHUP_baz.target
[Unit]
Description=Foo
[Service]
SendSIGHUP=baz

