# tests/auto/InvalidValue/2.31/bad/Service_PrivateDevices_foo.target
[Unit]
Description=Foo
[Service]
PrivateDevices=foo

