# tests/auto/InvalidValue/2.31/bad/Service_TTYVHangup_baz.target
[Unit]
Description=Foo
[Service]
TTYVHangup=baz

