# tests/auto/InvalidValue/2.31/bad/Service_BlockIOWriteBandwidthMax_--tmp--dev%2010000000000000T.target
[Unit]
Description=Foo
[Service]
BlockIOWriteBandwidthMax=/tmp/dev 10000000000000T

