# tests/auto/InvalidValue/2.31/bad/Service_Group_meisnotagroupcozofmyloooooooooooooooooooooooooooooooooognamegee.target
[Unit]
Description=Foo
[Service]
Group=meisnotagroupcozofmyloooooooooooooooooooooooooooooooooognamegee

