# tests/auto/InvalidValue/2.31/bad/Service_BlockIODeviceWeight_--tmp%2010001.target
[Unit]
Description=Foo
[Service]
BlockIODeviceWeight=/tmp 10001

