# tests/auto/InvalidValue/2.31/bad/Service_Environment_1.target
[Unit]
Description=Foo
[Service]
Environment=1

