# tests/auto/InvalidValue/2.31/bad/Service_MemoryDenyWriteExecute_foo.target
[Unit]
Description=Foo
[Service]
MemoryDenyWriteExecute=foo

