# tests/auto/InvalidValue/2.31/bad/Service_MemoryHigh_10000000000000K.target
[Unit]
Description=Foo
[Service]
MemoryHigh=10000000000000K

