# tests/auto/InvalidValue/2.31/bad/Service_MemoryAccounting_bar.target
[Unit]
Description=Foo
[Service]
MemoryAccounting=bar

