# tests/auto/InvalidValue/2.31/bad/Service_BlockIOReadBandwidthMax_--tmp--dev%2010000000000000T.target
[Unit]
Description=Foo
[Service]
BlockIOReadBandwidthMax=/tmp/dev 10000000000000T

