# tests/auto/InvalidValue/2.31/bad/Service_PrivateNetwork_bar.target
[Unit]
Description=Foo
[Service]
PrivateNetwork=bar

