# tests/auto/InvalidValue/2.31/bad/Service_MemorySwapMax_-1G.target
[Unit]
Description=Foo
[Service]
MemorySwapMax=-1G

