# tests/auto/InvalidValue/2.31/bad/Service_StartLimitIntervalSec_2potatoes.target
[Unit]
Description=Foo
[Service]
StartLimitIntervalSec=2potatoes

