# tests/auto/InvalidValue/2.31/bad/Service_SupplementaryGroups_4.target
[Unit]
Description=Foo
[Service]
SupplementaryGroups=4

