# tests/auto/InvalidValue/2.31/bad/Service_NonBlocking_foo.target
[Unit]
Description=Foo
[Service]
NonBlocking=foo

