# tests/auto/InvalidValue/2.31/bad/Service_RuntimeMaxSec_True.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=True

