# tests/auto/InvalidValue/2.31/bad/Service_TTYReset_bar.target
[Unit]
Description=Foo
[Service]
TTYReset=bar

