# tests/auto/InvalidValue/2.31/bad/Service_TimeoutStopSec_True.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=True

