# tests/auto/InvalidValue/2.31/bad/Service_TTYVTDisallocate_baz.target
[Unit]
Description=Foo
[Service]
TTYVTDisallocate=baz

