# tests/auto/InvalidValue/2.31/bad/Service_SyslogLevelPrefix_foo.target
[Unit]
Description=Foo
[Service]
SyslogLevelPrefix=foo

