# tests/auto/InvalidValue/2.31/bad/Service_LimitMEMLOCK_-1.target
[Unit]
Description=Foo
[Service]
LimitMEMLOCK=-1

