# tests/auto/InvalidValue/2.31/bad/Service_PrivateNetwork_foo.target
[Unit]
Description=Foo
[Service]
PrivateNetwork=foo

