# tests/auto/InvalidValue/2.31/bad/Service_TimeoutStopSec_1p.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=1p

