# tests/auto/InvalidValue/2.31/bad/Service_CPUQuota_10000000000000%25.target
[Unit]
Description=Foo
[Service]
CPUQuota=10000000000000%

