# tests/auto/InvalidValue/2.31/bad/Service_TTYReset_baz.target
[Unit]
Description=Foo
[Service]
TTYReset=baz

