# tests/auto/InvalidValue/2.31/bad/Service_TimerAccuracySec_2potatoes.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=2potatoes

