# tests/auto/InvalidValue/2.31/bad/Service_LimitDATA_-1.target
[Unit]
Description=Foo
[Service]
LimitDATA=-1

