# tests/auto/InvalidValue/2.31/bad/Service_BlockIOReadBandwidthMax_--tmp%20-1B.target
[Unit]
Description=Foo
[Service]
BlockIOReadBandwidthMax=/tmp -1B

