# tests/auto/InvalidValue/2.31/bad/Service_MemoryLow_10000000000000T.target
[Unit]
Description=Foo
[Service]
MemoryLow=10000000000000T

