# tests/auto/InvalidValue/2.31/bad/Service_CPUSchedulingResetOnFork_bar.target
[Unit]
Description=Foo
[Service]
CPUSchedulingResetOnFork=bar

