# tests/auto/InvalidValue/2.31/bad/Service_TTYVTDisallocate_foo.target
[Unit]
Description=Foo
[Service]
TTYVTDisallocate=foo

