# tests/auto/InvalidValue/2.31/bad/Service_RestrictRealtime_bar.target
[Unit]
Description=Foo
[Service]
RestrictRealtime=bar

