# tests/auto/InvalidValue/2.31/bad/Service_FileDescriptorStoreMax_-1.target
[Unit]
Description=Foo
[Service]
FileDescriptorStoreMax=-1

