# tests/auto/InvalidValue/2.31/bad/Service_BlockIODeviceWeight_--tmp%200.target
[Unit]
Description=Foo
[Service]
BlockIODeviceWeight=/tmp 0

