# tests/auto/InvalidValue/2.31/bad/Service_StartLimitBurst_-1.target
[Unit]
Description=Foo
[Service]
StartLimitBurst=-1

