# tests/auto/InvalidValue/2.31/bad/Service_PrivateDevices_bar.target
[Unit]
Description=Foo
[Service]
PrivateDevices=bar

