# tests/auto/InvalidValue/2.31/bad/Service_LimitSIGPENDING_-1.target
[Unit]
Description=Foo
[Service]
LimitSIGPENDING=-1

