# tests/auto/InvalidValue/2.31/bad/Service_PrivateDevices_baz.target
[Unit]
Description=Foo
[Service]
PrivateDevices=baz

