# tests/auto/InvalidValue/2.31/bad/Socket_SocketGroup_True.target
[Unit]
Description=Foo
[Socket]
SocketGroup=True

