# tests/auto/InvalidValue/2.31/bad/Socket_Mark_10000000000000.target
[Unit]
Description=Foo
[Socket]
Mark=10000000000000

