# tests/auto/InvalidValue/2.31/bad/Socket_PrivateDevices_baz.target
[Unit]
Description=Foo
[Socket]
PrivateDevices=baz

