# tests/auto/InvalidValue/2.31/bad/Socket_User_4.target
[Unit]
Description=Foo
[Socket]
User=4

