# tests/auto/InvalidValue/2.31/bad/Socket_TimeoutSec_0.5d.target
[Unit]
Description=Foo
[Socket]
TimeoutSec=0.5d

