# tests/auto/InvalidValue/2.31/bad/Socket_LimitNICE_10000000000000.target
[Unit]
Description=Foo
[Socket]
LimitNICE=10000000000000

