# tests/auto/InvalidValue/2.31/bad/Socket_LimitNOFILE_-1.target
[Unit]
Description=Foo
[Socket]
LimitNOFILE=-1

