# tests/auto/InvalidValue/2.31/bad/Socket_LimitAS_-1.target
[Unit]
Description=Foo
[Socket]
LimitAS=-1

