# tests/auto/InvalidValue/2.31/bad/Socket_SupplementaryGroups_3.target
[Unit]
Description=Foo
[Socket]
SupplementaryGroups=3

