# tests/auto/InvalidValue/2.31/bad/Socket_TasksAccounting_bar.target
[Unit]
Description=Foo
[Socket]
TasksAccounting=bar

