# tests/auto/InvalidValue/2.31/bad/Socket_Group_4.target
[Unit]
Description=Foo
[Socket]
Group=4

