# tests/auto/InvalidValue/2.31/bad/Socket_SupplementaryGroups_meisnotagroupcozofmyloooooooooooooooooooooooooooooooooognamegee.target
[Unit]
Description=Foo
[Socket]
SupplementaryGroups=meisnotagroupcozofmyloooooooooooooooooooooooooooooooooognamegee

