# tests/auto/InvalidValue/2.31/bad/Socket_ReusePort_bar.target
[Unit]
Description=Foo
[Socket]
ReusePort=bar

