# tests/auto/InvalidValue/2.31/bad/Socket_BindToDevice_dev--foo.target
[Unit]
Description=Foo
[Socket]
BindToDevice=dev/foo

