# tests/auto/InvalidValue/2.31/bad/Socket_ListenSpecial_1.target
[Unit]
Description=Foo
[Socket]
ListenSpecial=1

