# tests/auto/InvalidValue/2.31/bad/Socket_DeviceAllow_--tmp.target
[Unit]
Description=Foo
[Socket]
DeviceAllow=/tmp

