# tests/auto/InvalidValue/2.31/bad/Socket_IgnoreSIGPIPE_bar.target
[Unit]
Description=Foo
[Socket]
IgnoreSIGPIPE=bar

