# tests/auto/InvalidValue/2.31/bad/Socket_EnvironmentFile_tmp.target
[Unit]
Description=Foo
[Socket]
EnvironmentFile=tmp

