# tests/auto/InvalidValue/2.31/bad/Socket_StartupCPUWeight_0.target
[Unit]
Description=Foo
[Socket]
StartupCPUWeight=0

