# tests/auto/InvalidValue/2.31/bad/Socket_PassSecurity_bar.target
[Unit]
Description=Foo
[Socket]
PassSecurity=bar

