# tests/auto/InvalidValue/2.31/bad/Socket_SocketGroup_4.target
[Unit]
Description=Foo
[Socket]
SocketGroup=4

