# tests/auto/InvalidValue/2.31/bad/Socket_Group_meisnotagroupcozofmyloooooooooooooooooooooooooooooooooognamegee.target
[Unit]
Description=Foo
[Socket]
Group=meisnotagroupcozofmyloooooooooooooooooooooooooooooooooognamegee

