# tests/auto/InvalidValue/2.31/bad/Socket_EnvironmentFile_True.target
[Unit]
Description=Foo
[Socket]
EnvironmentFile=True

