# tests/auto/InvalidValue/2.31/bad/Socket_PipeSize_10000000000000.target
[Unit]
Description=Foo
[Socket]
PipeSize=10000000000000

