# tests/auto/InvalidValue/2.31/bad/Socket_Accept_baz.target
[Unit]
Description=Foo
[Socket]
Accept=baz

