# tests/auto/InvalidValue/2.31/bad/Socket_TTYPath_--devvvv.target
[Unit]
Description=Foo
[Socket]
TTYPath=/devvvv

