# tests/auto/InvalidValue/2.31/bad/Socket_LimitSTACK_-1.target
[Unit]
Description=Foo
[Socket]
LimitSTACK=-1

