# tests/auto/InvalidValue/2.31/bad/Socket_TTYVTDisallocate_bar.target
[Unit]
Description=Foo
[Socket]
TTYVTDisallocate=bar

