# tests/auto/InvalidValue/2.31/bad/Socket_BlockIODeviceWeight_--devvvv%2010001.target
[Unit]
Description=Foo
[Socket]
BlockIODeviceWeight=/devvvv 10001

