# tests/auto/InvalidValue/2.31/bad/Socket_RestrictRealtime_baz.target
[Unit]
Description=Foo
[Socket]
RestrictRealtime=baz

