# tests/auto/InvalidValue/2.31/bad/Socket_IgnoreSIGPIPE_baz.target
[Unit]
Description=Foo
[Socket]
IgnoreSIGPIPE=baz

