# tests/auto/InvalidValue/2.31/bad/Socket_CPUSchedulingResetOnFork_baz.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingResetOnFork=baz

