# tests/auto/InvalidValue/2.31/bad/Socket_LimitCPU_10000000000000.target
[Unit]
Description=Foo
[Socket]
LimitCPU=10000000000000

