# tests/auto/InvalidValue/2.31/bad/Socket_FreeBind_foo.target
[Unit]
Description=Foo
[Socket]
FreeBind=foo

