# tests/auto/InvalidValue/2.31/bad/Socket_PassCredentials_foo.target
[Unit]
Description=Foo
[Socket]
PassCredentials=foo

