# tests/auto/InvalidValue/2.31/bad/Socket_TTYReset_baz.target
[Unit]
Description=Foo
[Socket]
TTYReset=baz

