# tests/auto/InvalidValue/2.31/bad/Socket_NoDelay_baz.target
[Unit]
Description=Foo
[Socket]
NoDelay=baz

