# tests/auto/InvalidValue/2.31/bad/Socket_ReusePort_baz.target
[Unit]
Description=Foo
[Socket]
ReusePort=baz

