# tests/auto/InvalidValue/2.31/bad/Socket_BlockIOReadBandwidthMax_dev--foo%2010000000000000B.target
[Unit]
Description=Foo
[Socket]
BlockIOReadBandwidthMax=dev/foo 10000000000000B

