# tests/auto/InvalidValue/2.31/bad/Socket_LimitRTTIME_-1.target
[Unit]
Description=Foo
[Socket]
LimitRTTIME=-1

