# tests/auto/InvalidValue/2.31/bad/Socket_LimitMEMLOCK_10000000000000.target
[Unit]
Description=Foo
[Socket]
LimitMEMLOCK=10000000000000

