# tests/auto/InvalidValue/2.31/bad/Socket_SupplementaryGroups_True.target
[Unit]
Description=Foo
[Socket]
SupplementaryGroups=True

