# tests/auto/InvalidValue/2.31/bad/Socket_TTYVTDisallocate_foo.target
[Unit]
Description=Foo
[Socket]
TTYVTDisallocate=foo

