# tests/auto/InvalidValue/2.31/bad/Socket_SocketUser_True.target
[Unit]
Description=Foo
[Socket]
SocketUser=True

