# tests/auto/InvalidValue/2.31/bad/Socket_User_3.target
[Unit]
Description=Foo
[Socket]
User=3

