# tests/auto/InvalidValue/2.31/bad/Socket_TriggerLimitIntervalSec_0.5d.target
[Unit]
Description=Foo
[Socket]
TriggerLimitIntervalSec=0.5d

