# tests/auto/InvalidValue/2.31/bad/Socket_IOSchedulingPriority_8.target
[Unit]
Description=Foo
[Socket]
IOSchedulingPriority=8

