# tests/auto/InvalidValue/2.31/bad/Socket_RemoveIPC_baz.target
[Unit]
Description=Foo
[Socket]
RemoveIPC=baz

