# tests/auto/InvalidValue/2.31/bad/Socket_BlockIOWriteBandwidthMax_--devvvv%20-1B.target
[Unit]
Description=Foo
[Socket]
BlockIOWriteBandwidthMax=/devvvv -1B

