# tests/auto/InvalidValue/2.31/bad/Socket_CPUSchedulingResetOnFork_bar.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingResetOnFork=bar

