# tests/auto/InvalidValue/2.31/bad/Socket_TriggerLimitBurst_True.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=True

