# tests/auto/InvalidValue/2.31/bad/Socket_LimitRSS_10000000000000.target
[Unit]
Description=Foo
[Socket]
LimitRSS=10000000000000

