# tests/auto/InvalidValue/2.31/bad/Socket_BindToDevice_--tmp.target
[Unit]
Description=Foo
[Socket]
BindToDevice=/tmp

