# tests/auto/InvalidValue/2.31/bad/Socket_BlockIODeviceWeight_--tmp--dev%2010001.target
[Unit]
Description=Foo
[Socket]
BlockIODeviceWeight=/tmp/dev 10001

