# tests/auto/InvalidValue/2.31/bad/Socket_MemoryLow_-1T.target
[Unit]
Description=Foo
[Socket]
MemoryLow=-1T

