# tests/auto/InvalidValue/2.31/bad/Socket_IPTOS_10000000000000.target
[Unit]
Description=Foo
[Socket]
IPTOS=10000000000000

