# tests/auto/InvalidValue/2.31/bad/Socket_SupplementaryGroups_4.target
[Unit]
Description=Foo
[Socket]
SupplementaryGroups=4

