# tests/auto/InvalidValue/2.31/bad/Socket_RestrictRealtime_foo.target
[Unit]
Description=Foo
[Socket]
RestrictRealtime=foo

