# tests/auto/InvalidValue/2.31/bad/Socket_MessageQueueMessageSize_10000000000000.target
[Unit]
Description=Foo
[Socket]
MessageQueueMessageSize=10000000000000

