# tests/auto/InvalidValue/2.31/bad/Socket_PipeSize_-1.target
[Unit]
Description=Foo
[Socket]
PipeSize=-1

