# tests/auto/InvalidValue/2.31/bad/Socket_BlockStartupIOWeight_0.target
[Unit]
Description=Foo
[Socket]
BlockStartupIOWeight=0

