# tests/auto/InvalidValue/2.31/bad/Socket_TimeoutSec_True.target
[Unit]
Description=Foo
[Socket]
TimeoutSec=True

