# tests/auto/InvalidValue/2.31/bad/Socket_BlockIOReadBandwidthMax_--tmp%20-1B.target
[Unit]
Description=Foo
[Socket]
BlockIOReadBandwidthMax=/tmp -1B

