# tests/auto/InvalidValue/2.31/bad/Socket_MemoryHigh_10000000000000K.target
[Unit]
Description=Foo
[Socket]
MemoryHigh=10000000000000K

