# tests/auto/InvalidValue/2.31/bad/Socket_LimitFSIZE_10000000000000.target
[Unit]
Description=Foo
[Socket]
LimitFSIZE=10000000000000

