# tests/auto/InvalidValue/2.31/bad/Socket_RemoveOnStop_bar.target
[Unit]
Description=Foo
[Socket]
RemoveOnStop=bar

