# tests/auto/InvalidValue/2.31/bad/Socket_BlockIOWeight_10001.target
[Unit]
Description=Foo
[Socket]
BlockIOWeight=10001

