# tests/auto/InvalidValue/2.31/bad/Socket_LimitAS_10000000000000.target
[Unit]
Description=Foo
[Socket]
LimitAS=10000000000000

