# tests/auto/InvalidValue/2.31/bad/Socket_KeepAlive_baz.target
[Unit]
Description=Foo
[Socket]
KeepAlive=baz

