# tests/auto/InvalidValue/2.31/bad/Socket_KeepAlive_foo.target
[Unit]
Description=Foo
[Socket]
KeepAlive=foo

