# tests/auto/InvalidValue/2.31/bad/Socket_IOSchedulingPriority_-1.target
[Unit]
Description=Foo
[Socket]
IOSchedulingPriority=-1

