# tests/auto/InvalidValue/2.31/bad/Socket_RestrictSUIDSGID_foo.target
[Unit]
Description=Foo
[Socket]
RestrictSUIDSGID=foo

