# tests/auto/InvalidValue/2.31/bad/Socket_BlockIOWriteBandwidthMax_--tmp--dev%20-1M.target
[Unit]
Description=Foo
[Socket]
BlockIOWriteBandwidthMax=/tmp/dev -1M

