# tests/auto/InvalidValue/2.31/bad/Socket_NoDelay_bar.target
[Unit]
Description=Foo
[Socket]
NoDelay=bar

