# tests/auto/InvalidValue/2.31/bad/Socket_LimitDATA_10000000000000.target
[Unit]
Description=Foo
[Socket]
LimitDATA=10000000000000

