# tests/auto/InvalidValue/2.31/bad/Socket_RemoveOnStop_baz.target
[Unit]
Description=Foo
[Socket]
RemoveOnStop=baz

