# tests/auto/InvalidValue/2.31/bad/Socket_PassCredentials_bar.target
[Unit]
Description=Foo
[Socket]
PassCredentials=bar

