# tests/auto/InvalidValue/2.31/bad/Socket_PrivateTmp_bar.target
[Unit]
Description=Foo
[Socket]
PrivateTmp=bar

