# tests/auto/InvalidValue/2.31/bad/Socket_BlockIODeviceWeight_--tmp%200.target
[Unit]
Description=Foo
[Socket]
BlockIODeviceWeight=/tmp 0

