# tests/auto/InvalidValue/2.31/bad/Socket_MemoryDenyWriteExecute_foo.target
[Unit]
Description=Foo
[Socket]
MemoryDenyWriteExecute=foo

