# tests/auto/InvalidValue/2.31/bad/Socket_SocketMode_999.target
[Unit]
Description=Foo
[Socket]
SocketMode=999

