# tests/auto/InvalidValue/2.31/bad/Socket_PassSecurity_baz.target
[Unit]
Description=Foo
[Socket]
PassSecurity=baz

