# tests/auto/InvalidValue/2.31/bad/Socket_PrivateNetwork_baz.target
[Unit]
Description=Foo
[Socket]
PrivateNetwork=baz

