# tests/auto/InvalidValue/2.31/bad/Socket_MessageQueueMaxMessages_10000000000000.target
[Unit]
Description=Foo
[Socket]
MessageQueueMaxMessages=10000000000000

