# tests/auto/InvalidValue/2.31/bad/Socket_PassCredentials_baz.target
[Unit]
Description=Foo
[Socket]
PassCredentials=baz

