# tests/auto/InvalidValue/2.31/bad/Socket_LimitSIGPENDING_10000000000000.target
[Unit]
Description=Foo
[Socket]
LimitSIGPENDING=10000000000000

