# tests/auto/InvalidValue/2.31/bad/Socket_LimitMEMLOCK_-1.target
[Unit]
Description=Foo
[Socket]
LimitMEMLOCK=-1

