# tests/auto/InvalidValue/2.31/bad/Socket_MemoryLow_10000000000000B.target
[Unit]
Description=Foo
[Socket]
MemoryLow=10000000000000B

