# tests/auto/InvalidValue/2.31/bad/Socket_MemoryLimit_-1G.target
[Unit]
Description=Foo
[Socket]
MemoryLimit=-1G

