# tests/auto/InvalidValue/2.31/bad/Socket_Group_3.target
[Unit]
Description=Foo
[Socket]
Group=3

