# tests/auto/InvalidValue/2.31/bad/Socket_TTYVTDisallocate_baz.target
[Unit]
Description=Foo
[Socket]
TTYVTDisallocate=baz

