# tests/auto/InvalidValue/2.31/bad/Socket_BlockIOAccounting_bar.target
[Unit]
Description=Foo
[Socket]
BlockIOAccounting=bar

