# tests/auto/InvalidValue/2.31/bad/Socket_BindToDevice_--devvvv.target
[Unit]
Description=Foo
[Socket]
BindToDevice=/devvvv

