# tests/auto/InvalidValue/2.31/bad/Socket_StartupCPUShares_1.target
[Unit]
Description=Foo
[Socket]
StartupCPUShares=1

