# tests/auto/InvalidValue/2.31/bad/Socket_LimitFSIZE_-1.target
[Unit]
Description=Foo
[Socket]
LimitFSIZE=-1

