# tests/auto/InvalidValue/2.31/bad/Socket_LimitSTACK_10000000000000.target
[Unit]
Description=Foo
[Socket]
LimitSTACK=10000000000000

