# tests/auto/InvalidValue/2.31/bad/Socket_TasksAccounting_baz.target
[Unit]
Description=Foo
[Socket]
TasksAccounting=baz

