# tests/auto/InvalidValue/2.31/bad/Socket_Environment_True.target
[Unit]
Description=Foo
[Socket]
Environment=True

