# tests/auto/InvalidValue/2.31/bad/Socket_BlockIOReadBandwidthMax_--tmp%2010000000000000K.target
[Unit]
Description=Foo
[Socket]
BlockIOReadBandwidthMax=/tmp 10000000000000K

