# tests/auto/InvalidValue/2.31/bad/Socket_DeviceAllow_dev--foo.target
[Unit]
Description=Foo
[Socket]
DeviceAllow=dev/foo

