# tests/auto/InvalidValue/2.31/bad/Socket_PrivateNetwork_bar.target
[Unit]
Description=Foo
[Socket]
PrivateNetwork=bar

