# tests/auto/InvalidValue/2.31/bad/Socket_MemoryMax_-1M.target
[Unit]
Description=Foo
[Socket]
MemoryMax=-1M

