# tests/auto/InvalidValue/2.31/bad/Socket_ListenFIFO_..--..--.target
[Unit]
Description=Foo
[Socket]
ListenFIFO=../../

