# tests/auto/InvalidValue/2.31/bad/Socket_LimitMSGQUEUE_10000000000000.target
[Unit]
Description=Foo
[Socket]
LimitMSGQUEUE=10000000000000

