# tests/auto/InvalidValue/2.31/bad/Socket_PrivateTmp_foo.target
[Unit]
Description=Foo
[Socket]
PrivateTmp=foo

