# tests/auto/InvalidValue/2.31/bad/Socket_SELinuxContextFromNet_foo.target
[Unit]
Description=Foo
[Socket]
SELinuxContextFromNet=foo

