# tests/auto/InvalidValue/2.31/bad/Socket_LimitMSGQUEUE_-1.target
[Unit]
Description=Foo
[Socket]
LimitMSGQUEUE=-1

