# tests/auto/InvalidValue/2.31/bad/Socket_LimitSIGPENDING_-1.target
[Unit]
Description=Foo
[Socket]
LimitSIGPENDING=-1

