# tests/auto/InvalidValue/2.31/bad/Socket_Transparent_foo.target
[Unit]
Description=Foo
[Socket]
Transparent=foo

