# tests/auto/InvalidValue/2.31/bad/Socket_Symlinks_True.target
[Unit]
Description=Foo
[Socket]
Symlinks=True

