# tests/auto/InvalidValue/2.31/bad/Socket_Broadcast_bar.target
[Unit]
Description=Foo
[Socket]
Broadcast=bar

