# tests/auto/InvalidValue/2.31/bad/Socket_SocketMode_True.target
[Unit]
Description=Foo
[Socket]
SocketMode=True

