# tests/auto/InvalidValue/2.31/bad/Socket_MountAPIVFS_bar.target
[Unit]
Description=Foo
[Socket]
MountAPIVFS=bar

