# tests/auto/InvalidValue/2.31/bad/Socket_BlockIOReadBandwidthMax_--devvvv%20-1K.target
[Unit]
Description=Foo
[Socket]
BlockIOReadBandwidthMax=/devvvv -1K

