# tests/auto/InvalidValue/2.31/bad/Socket_CPUSchedulingResetOnFork_foo.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingResetOnFork=foo

