# tests/auto/InvalidValue/2.31/bad/Socket_SocketUser_4.target
[Unit]
Description=Foo
[Socket]
SocketUser=4

