# tests/auto/InvalidValue/2.31/bad/Socket_TimerSlackNSec_-1.target
[Unit]
Description=Foo
[Socket]
TimerSlackNSec=-1

