# tests/auto/InvalidValue/2.31/bad/Socket_BlockIOReadBandwidthMax_dev--foo%20-1G.target
[Unit]
Description=Foo
[Socket]
BlockIOReadBandwidthMax=dev/foo -1G

