# tests/auto/InvalidValue/2.31/bad/Swap_EnvironmentFile_tmp.target
[Unit]
Description=Foo
[Swap]
EnvironmentFile=tmp

