# tests/auto/InvalidValue/2.31/bad/Swap_OOMScoreAdjust_1001.target
[Unit]
Description=Foo
[Swap]
OOMScoreAdjust=1001

