# tests/auto/InvalidValue/2.31/bad/Swap_MemoryMax_10000000000000T.target
[Unit]
Description=Foo
[Swap]
MemoryMax=10000000000000T

