# tests/auto/InvalidValue/2.31/bad/Swap_BlockIOReadBandwidthMax_--tmp%20-1G.target
[Unit]
Description=Foo
[Swap]
BlockIOReadBandwidthMax=/tmp -1G

