# tests/auto/InvalidValue/2.31/bad/Swap_BlockIOWriteBandwidthMax_--tmp%2010000000000000B.target
[Unit]
Description=Foo
[Swap]
BlockIOWriteBandwidthMax=/tmp 10000000000000B

