# tests/auto/InvalidValue/2.31/bad/Swap_BlockIOReadBandwidthMax_dev--foo%20-1B.target
[Unit]
Description=Foo
[Swap]
BlockIOReadBandwidthMax=dev/foo -1B

