# tests/auto/InvalidValue/2.31/bad/Swap_RestrictRealtime_foo.target
[Unit]
Description=Foo
[Swap]
RestrictRealtime=foo

