# tests/auto/InvalidValue/2.31/bad/Swap_TasksAccounting_bar.target
[Unit]
Description=Foo
[Swap]
TasksAccounting=bar

