# tests/auto/InvalidValue/2.31/bad/Swap_PrivateTmp_baz.target
[Unit]
Description=Foo
[Swap]
PrivateTmp=baz

