# tests/auto/InvalidValue/2.31/bad/Swap_LimitMSGQUEUE_-1.target
[Unit]
Description=Foo
[Swap]
LimitMSGQUEUE=-1

