# tests/auto/InvalidValue/2.31/bad/Swap_StartupCPUWeight_0.target
[Unit]
Description=Foo
[Swap]
StartupCPUWeight=0

