# tests/auto/InvalidValue/2.31/bad/Swap_User_True.target
[Unit]
Description=Foo
[Swap]
User=True

