# tests/auto/InvalidValue/2.31/bad/Swap_LimitRTTIME_-1.target
[Unit]
Description=Foo
[Swap]
LimitRTTIME=-1

