# tests/auto/InvalidValue/2.31/bad/Swap_BlockStartupIOWeight_0.target
[Unit]
Description=Foo
[Swap]
BlockStartupIOWeight=0

