# tests/auto/InvalidValue/2.31/bad/Swap_BlockIODeviceWeight_--devvvv%2010001.target
[Unit]
Description=Foo
[Swap]
BlockIODeviceWeight=/devvvv 10001

