# tests/auto/InvalidValue/2.31/bad/Swap_SupplementaryGroups_True.target
[Unit]
Description=Foo
[Swap]
SupplementaryGroups=True

