# tests/auto/InvalidValue/2.31/bad/Swap_BlockIOWriteBandwidthMax_dev--foo%2010000000000000M.target
[Unit]
Description=Foo
[Swap]
BlockIOWriteBandwidthMax=dev/foo 10000000000000M

