# tests/auto/InvalidValue/2.31/bad/Swap_LimitAS_10000000000000.target
[Unit]
Description=Foo
[Swap]
LimitAS=10000000000000

