# tests/auto/InvalidValue/2.31/bad/Swap_LimitSTACK_10000000000000.target
[Unit]
Description=Foo
[Swap]
LimitSTACK=10000000000000

