# tests/auto/InvalidValue/2.31/bad/Swap_SyslogLevelPrefix_baz.target
[Unit]
Description=Foo
[Swap]
SyslogLevelPrefix=baz

