# tests/auto/InvalidValue/2.31/bad/Swap_IgnoreSIGPIPE_bar.target
[Unit]
Description=Foo
[Swap]
IgnoreSIGPIPE=bar

