# tests/auto/InvalidValue/2.31/bad/Swap_LimitNPROC_10000000000000.target
[Unit]
Description=Foo
[Swap]
LimitNPROC=10000000000000

