# tests/auto/InvalidValue/2.31/bad/Swap_MemoryHigh_-1T.target
[Unit]
Description=Foo
[Swap]
MemoryHigh=-1T

