# tests/auto/InvalidValue/2.31/bad/Swap_TasksAccounting_foo.target
[Unit]
Description=Foo
[Swap]
TasksAccounting=foo

