# tests/auto/InvalidValue/2.31/bad/Swap_CPUSchedulingPriority_100.target
[Unit]
Description=Foo
[Swap]
CPUSchedulingPriority=100

