# tests/auto/InvalidValue/2.31/bad/Swap_LimitSIGPENDING_-1.target
[Unit]
Description=Foo
[Swap]
LimitSIGPENDING=-1

