# tests/auto/InvalidValue/2.31/bad/Swap_LimitNOFILE_10000000000000.target
[Unit]
Description=Foo
[Swap]
LimitNOFILE=10000000000000

