# tests/auto/InvalidValue/2.31/bad/Swap_LimitRTPRIO_10000000000000.target
[Unit]
Description=Foo
[Swap]
LimitRTPRIO=10000000000000

