# tests/auto/InvalidValue/2.31/bad/Swap_LimitMEMLOCK_10000000000000.target
[Unit]
Description=Foo
[Swap]
LimitMEMLOCK=10000000000000

