# tests/auto/InvalidValue/2.31/bad/Swap_LimitLOCKS_-1.target
[Unit]
Description=Foo
[Swap]
LimitLOCKS=-1

