# tests/auto/InvalidValue/2.31/bad/Swap_TTYVTDisallocate_foo.target
[Unit]
Description=Foo
[Swap]
TTYVTDisallocate=foo

