# tests/auto/InvalidValue/2.31/bad/Swap_TimeoutSec_True.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=True

