# tests/auto/InvalidValue/2.31/bad/Swap_PrivateTmp_bar.target
[Unit]
Description=Foo
[Swap]
PrivateTmp=bar

