# tests/auto/InvalidValue/2.31/bad/Swap_IgnoreSIGPIPE_foo.target
[Unit]
Description=Foo
[Swap]
IgnoreSIGPIPE=foo

