# tests/auto/InvalidValue/2.31/bad/Swap_RestrictRealtime_bar.target
[Unit]
Description=Foo
[Swap]
RestrictRealtime=bar

