# tests/auto/InvalidValue/2.31/bad/Swap_MemoryAccounting_bar.target
[Unit]
Description=Foo
[Swap]
MemoryAccounting=bar

