# tests/auto/InvalidValue/2.31/bad/Swap_TimerSlackNSec_-1.target
[Unit]
Description=Foo
[Swap]
TimerSlackNSec=-1

