# tests/auto/InvalidValue/2.31/bad/Swap_MemoryAccounting_baz.target
[Unit]
Description=Foo
[Swap]
MemoryAccounting=baz

