# tests/auto/InvalidValue/2.31/bad/Swap_SyslogLevelPrefix_foo.target
[Unit]
Description=Foo
[Swap]
SyslogLevelPrefix=foo

