# tests/auto/InvalidValue/2.31/bad/Swap_UMask_10000000000000.target
[Unit]
Description=Foo
[Swap]
UMask=10000000000000

