# tests/auto/InvalidValue/2.31/bad/Swap_PrivateNetwork_foo.target
[Unit]
Description=Foo
[Swap]
PrivateNetwork=foo

