# tests/auto/InvalidValue/2.31/bad/Swap_LimitNOFILE_-1.target
[Unit]
Description=Foo
[Swap]
LimitNOFILE=-1

