# tests/auto/InvalidValue/2.31/bad/Swap_PrivateDevices_bar.target
[Unit]
Description=Foo
[Swap]
PrivateDevices=bar

