# tests/auto/InvalidValue/2.31/bad/Swap_LimitRTPRIO_-1.target
[Unit]
Description=Foo
[Swap]
LimitRTPRIO=-1

