# tests/auto/InvalidValue/2.31/bad/Swap_LimitNICE_-1.target
[Unit]
Description=Foo
[Swap]
LimitNICE=-1

