# tests/auto/InvalidValue/2.31/bad/Swap_BlockIOAccounting_bar.target
[Unit]
Description=Foo
[Swap]
BlockIOAccounting=bar

