# tests/auto/InvalidValue/2.31/bad/Swap_LimitAS_-1.target
[Unit]
Description=Foo
[Swap]
LimitAS=-1

