# tests/auto/InvalidValue/2.31/bad/Swap_IOAccounting_foo.target
[Unit]
Description=Foo
[Swap]
IOAccounting=foo

