# tests/auto/InvalidValue/2.31/bad/Swap_StartupIOWeight_10001.target
[Unit]
Description=Foo
[Swap]
StartupIOWeight=10001

