# tests/auto/InvalidValue/2.31/bad/Swap_RestrictRealtime_baz.target
[Unit]
Description=Foo
[Swap]
RestrictRealtime=baz

