# tests/auto/InvalidValue/2.31/bad/Swap_StartupCPUWeight_10001.target
[Unit]
Description=Foo
[Swap]
StartupCPUWeight=10001

