# tests/auto/InvalidValue/2.31/bad/Swap_RestrictSUIDSGID_baz.target
[Unit]
Description=Foo
[Swap]
RestrictSUIDSGID=baz

