# tests/auto/InvalidValue/2.31/bad/Swap_Group_meisnotagroupcozofmyloooooooooooooooooooooooooooooooooognamegee.target
[Unit]
Description=Foo
[Swap]
Group=meisnotagroupcozofmyloooooooooooooooooooooooooooooooooognamegee

