# tests/auto/InvalidValue/2.31/bad/Swap_SendSIGHUP_bar.target
[Unit]
Description=Foo
[Swap]
SendSIGHUP=bar

