# tests/auto/InvalidValue/2.31/bad/Swap_TTYPath_--devvvv.target
[Unit]
Description=Foo
[Swap]
TTYPath=/devvvv

