# tests/auto/InvalidValue/2.31/bad/Swap_MemoryLimit_10000000000000T.target
[Unit]
Description=Foo
[Swap]
MemoryLimit=10000000000000T

