# tests/auto/InvalidValue/2.31/bad/Swap_CPUQuota_-1%25.target
[Unit]
Description=Foo
[Swap]
CPUQuota=-1%

