# tests/auto/InvalidValue/2.31/bad/Swap_DeviceAllow_dev--foo.target
[Unit]
Description=Foo
[Swap]
DeviceAllow=dev/foo

