# tests/auto/InvalidValue/2.31/bad/Swap_LimitFSIZE_10000000000000.target
[Unit]
Description=Foo
[Swap]
LimitFSIZE=10000000000000

