# tests/auto/InvalidValue/2.31/bad/Swap_LimitDATA_10000000000000.target
[Unit]
Description=Foo
[Swap]
LimitDATA=10000000000000

