# tests/auto/InvalidValue/2.31/bad/Swap_MountAPIVFS_foo.target
[Unit]
Description=Foo
[Swap]
MountAPIVFS=foo

