# tests/auto/InvalidValue/2.31/bad/Swap_RemoveIPC_baz.target
[Unit]
Description=Foo
[Swap]
RemoveIPC=baz

