# tests/auto/InvalidValue/2.31/bad/Swap_CPUShares_1.target
[Unit]
Description=Foo
[Swap]
CPUShares=1

