# tests/auto/InvalidValue/2.31/bad/Swap_BlockIODeviceWeight_--tmp--dev%200.target
[Unit]
Description=Foo
[Swap]
BlockIODeviceWeight=/tmp/dev 0

