# tests/auto/InvalidValue/2.31/bad/Swap_BlockIOAccounting_foo.target
[Unit]
Description=Foo
[Swap]
BlockIOAccounting=foo

