# tests/auto/InvalidValue/2.31/bad/Swap_MemoryHigh_10000000000000B.target
[Unit]
Description=Foo
[Swap]
MemoryHigh=10000000000000B

