# tests/auto/InvalidValue/2.31/bad/Swap_SendSIGKILL_bar.target
[Unit]
Description=Foo
[Swap]
SendSIGKILL=bar

