# tests/auto/InvalidValue/2.31/bad/Swap_LimitSIGPENDING_10000000000000.target
[Unit]
Description=Foo
[Swap]
LimitSIGPENDING=10000000000000

