# tests/auto/InvalidValue/2.31/bad/Swap_IOSchedulingPriority_-1.target
[Unit]
Description=Foo
[Swap]
IOSchedulingPriority=-1

