# tests/auto/InvalidValue/2.31/bad/Swap_SupplementaryGroups_meisnotagroupcozofmyloooooooooooooooooooooooooooooooooognamegee.target
[Unit]
Description=Foo
[Swap]
SupplementaryGroups=meisnotagroupcozofmyloooooooooooooooooooooooooooooooooognamegee

