# tests/auto/InvalidValue/2.31/bad/Swap_MemoryLow_-1T.target
[Unit]
Description=Foo
[Swap]
MemoryLow=-1T

