# tests/auto/InvalidValue/2.31/bad/Swap_TTYReset_bar.target
[Unit]
Description=Foo
[Swap]
TTYReset=bar

