# tests/auto/InvalidValue/2.31/bad/Swap_PrivateDevices_baz.target
[Unit]
Description=Foo
[Swap]
PrivateDevices=baz

