# tests/auto/InvalidValue/2.31/bad/Swap_LimitRSS_-1.target
[Unit]
Description=Foo
[Swap]
LimitRSS=-1

