# tests/auto/InvalidValue/2.31/bad/Swap_TasksMax_10000000000000.target
[Unit]
Description=Foo
[Swap]
TasksMax=10000000000000

