# tests/auto/InvalidValue/2.31/bad/Tap_OneQueue_baz.target
[Unit]
Description=Foo
[Tap]
OneQueue=baz

