# tests/auto/InvalidValue/2.31/bad/Tap_MultiQueue_baz.target
[Unit]
Description=Foo
[Tap]
MultiQueue=baz

