# tests/auto/InvalidValue/2.31/bad/Tap_PacketInfo_bar.target
[Unit]
Description=Foo
[Tap]
PacketInfo=bar

