# tests/auto/InvalidValue/2.31/bad/Tap_PacketInfo_foo.target
[Unit]
Description=Foo
[Tap]
PacketInfo=foo

