# tests/auto/InvalidValue/2.31/bad/Tap_VNetHeader_baz.target
[Unit]
Description=Foo
[Tap]
VNetHeader=baz

