# tests/auto/InvalidValue/2.31/bad/Tap_User_meisnotausercozofmyloooooooooooooooooooooooooooooooooognamegee.target
[Unit]
Description=Foo
[Tap]
User=meisnotausercozofmyloooooooooooooooooooooooooooooooooognamegee

