# tests/auto/InvalidValue/2.31/bad/Timer_OnUnitActiveSec_3decades.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=3decades

