# tests/auto/InvalidValue/2.31/bad/Timer_OnBootSec_1p.target
[Unit]
Description=Foo
[Timer]
OnBootSec=1p

