# tests/auto/InvalidValue/2.31/bad/Timer_AccuracySec_0.5d.target
[Unit]
Description=Foo
[Timer]
AccuracySec=0.5d

