# tests/auto/InvalidValue/2.31/bad/Timer_OnActiveSec_True.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=True

