# tests/auto/InvalidValue/2.31/bad/Timer_OnClockChange_baz.target
[Unit]
Description=Foo
[Timer]
OnClockChange=baz

