# tests/auto/InvalidValue/2.31/bad/Timer_OnUnitInactiveSec_3decades.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=3decades

