# tests/auto/InvalidValue/2.31/bad/Timer_OnActiveSec_0.5d.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=0.5d

