# tests/auto/InvalidValue/2.31/bad/Timer_OnTimezoneChange_bar.target
[Unit]
Description=Foo
[Timer]
OnTimezoneChange=bar

