# tests/auto/InvalidValue/2.31/bad/Time_PollIntervalMaxSec_16.target
[Unit]
Description=Foo
[Time]
PollIntervalMaxSec=16

