# tests/auto/InvalidValue/2.31/bad/Timer_OnStartupSec_True.target
[Unit]
Description=Foo
[Timer]
OnStartupSec=True

