# tests/auto/InvalidValue/2.31/bad/Timer_OnUnitActiveSec_1p.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=1p

