# tests/auto/InvalidValue/2.31/bad/Timer_WakeSystem_foo.target
[Unit]
Description=Foo
[Timer]
WakeSystem=foo

