# tests/auto/InvalidValue/2.31/bad/Timer_RandomizedDelaySec_True.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=True

