# tests/auto/InvalidValue/2.31/bad/Timer_OnBootSec_0.5d.target
[Unit]
Description=Foo
[Timer]
OnBootSec=0.5d

