# tests/auto/InvalidValue/2.31/bad/Timer_OnClockChange_foo.target
[Unit]
Description=Foo
[Timer]
OnClockChange=foo

