# tests/auto/InvalidValue/2.31/bad/Timer_Persistent_foo.target
[Unit]
Description=Foo
[Timer]
Persistent=foo

