# tests/auto/InvalidValue/2.31/bad/Timer_Persistent_bar.target
[Unit]
Description=Foo
[Timer]
Persistent=bar

