# tests/auto/InvalidValue/2.31/bad/Timer_OnUnitInactiveSec_True.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=True

