# tests/auto/InvalidValue/2.31/bad/Timer_RandomizedDelaySec_0.5d.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=0.5d

