# tests/auto/InvalidValue/2.31/bad/Tun_MultiQueue_baz.target
[Unit]
Description=Foo
[Tun]
MultiQueue=baz

