# tests/auto/InvalidValue/2.31/bad/Tun_OneQueue_baz.target
[Unit]
Description=Foo
[Tun]
OneQueue=baz

