# tests/auto/InvalidValue/2.31/bad/Tun_Group_True.target
[Unit]
Description=Foo
[Tun]
Group=True

