# tests/auto/InvalidValue/2.31/bad/Tun_OneQueue_foo.target
[Unit]
Description=Foo
[Tun]
OneQueue=foo

