# tests/auto/InvalidValue/2.31/bad/Tun_User_meisnotausercozofmyloooooooooooooooooooooooooooooooooognamegee.target
[Unit]
Description=Foo
[Tun]
User=meisnotausercozofmyloooooooooooooooooooooooooooooooooognamegee

