# tests/auto/InvalidValue/2.31/bad/Tun_MultiQueue_bar.target
[Unit]
Description=Foo
[Tun]
MultiQueue=bar

