# tests/auto/InvalidValue/2.31/bad/Tunnel_CopyDSCP_foo.target
[Unit]
Description=Foo
[Tunnel]
CopyDSCP=foo

