# tests/auto/InvalidValue/2.31/bad/Tunnel_ISATAP_foo.target
[Unit]
Description=Foo
[Tunnel]
ISATAP=foo

