# tests/auto/InvalidValue/2.31/bad/Tunnel_Independent_foo.target
[Unit]
Description=Foo
[Tunnel]
Independent=foo

