# tests/auto/InvalidValue/2.31/bad/Tunnel_CopyDSCP_baz.target
[Unit]
Description=Foo
[Tunnel]
CopyDSCP=baz

