# tests/auto/InvalidValue/2.31/bad/Tunnel_TTL_256.target
[Unit]
Description=Foo
[Tunnel]
TTL=256

