# tests/auto/InvalidValue/2.31/bad/Tunnel_ERSPANIndex_0.target
[Unit]
Description=Foo
[Tunnel]
ERSPANIndex=0

