# tests/auto/InvalidValue/2.31/bad/Tunnel_Independent_bar.target
[Unit]
Description=Foo
[Tunnel]
Independent=bar

