# tests/auto/InvalidValue/2.31/bad/Tunnel_FooOverUDP_baz.target
[Unit]
Description=Foo
[Tunnel]
FooOverUDP=baz

