# tests/auto/InvalidValue/2.31/bad/Tunnel_ISATAP_bar.target
[Unit]
Description=Foo
[Tunnel]
ISATAP=bar

