# tests/auto/InvalidValue/2.31/bad/Tunnel_SerializeTunneledPackets_foo.target
[Unit]
Description=Foo
[Tunnel]
SerializeTunneledPackets=foo

