# tests/auto/InvalidValue/2.31/bad/Tunnel_EncapsulationLimit_-1.target
[Unit]
Description=Foo
[Tunnel]
EncapsulationLimit=-1

