# tests/auto/InvalidValue/2.31/bad/Tunnel_DiscoverPathMTU_bar.target
[Unit]
Description=Foo
[Tunnel]
DiscoverPathMTU=bar

