# tests/auto/InvalidValue/2.31/bad/Tunnel_FOUDestinationPort_0.target
[Unit]
Description=Foo
[Tunnel]
FOUDestinationPort=0

