# tests/auto/InvalidValue/2.31/bad/Tunnel_FOUDestinationPort_65535.target
[Unit]
Description=Foo
[Tunnel]
FOUDestinationPort=65535

