# tests/auto/InvalidValue/2.31/bad/Tunnel_EncapsulationLimit_256.target
[Unit]
Description=Foo
[Tunnel]
EncapsulationLimit=256

