# tests/auto/InvalidValue/2.31/bad/Tunnel_ISATAP_baz.target
[Unit]
Description=Foo
[Tunnel]
ISATAP=baz

