# tests/auto/InvalidValue/2.31/bad/Unit_ConditionFileIsExecutable_1.target
[Unit]
Description=Foo
Description=abc
ConditionFileIsExecutable=1

