# tests/auto/InvalidValue/2.31/bad/Unit_JobRunningTimeoutSec_True.target
[Unit]
Description=Foo
Description=abc
JobRunningTimeoutSec=True

