# tests/auto/InvalidValue/2.31/bad/Unit_AssertPathExists_1.target
[Unit]
Description=Foo
Description=abc
AssertPathExists=1

