# tests/auto/InvalidValue/2.31/bad/Unit_AssertFileIsExecutable_1.target
[Unit]
Description=Foo
Description=abc
AssertFileIsExecutable=1

