# tests/auto/InvalidValue/2.31/bad/Unit_AssertCPUs_-1.target
[Unit]
Description=Foo
Description=abc
AssertCPUs=-1

