# tests/auto/InvalidValue/2.31/bad/Unit_JobTimeoutSec_True.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=True

