# tests/auto/InvalidValue/2.31/bad/Unit_DefaultDependencies_baz.target
[Unit]
Description=Foo
Description=abc
DefaultDependencies=baz

