# tests/auto/InvalidValue/2.31/bad/Unit_StartLimitBurst_10000000000000.target
[Unit]
Description=Foo
Description=abc
StartLimitBurst=10000000000000

