# tests/auto/InvalidValue/2.31/bad/Unit_StartLimitBurst_-1.target
[Unit]
Description=Foo
Description=abc
StartLimitBurst=-1

