# tests/auto/InvalidValue/2.31/bad/Unit_RefuseManualStop_bar.target
[Unit]
Description=Foo
Description=abc
RefuseManualStop=bar

