# tests/auto/InvalidValue/2.31/bad/Unit_ConditionPathIsReadWrite_..--..--.target
[Unit]
Description=Foo
Description=abc
ConditionPathIsReadWrite=../../

