# tests/auto/InvalidValue/2.31/bad/Unit_ConditionFileIsExecutable_True.target
[Unit]
Description=Foo
Description=abc
ConditionFileIsExecutable=True

