# tests/auto/InvalidValue/2.31/bad/Unit_IgnoreOnIsolate_bar.target
[Unit]
Description=Foo
Description=abc
IgnoreOnIsolate=bar

