# tests/auto/InvalidValue/2.31/bad/Unit_AssertFileNotEmpty_1.target
[Unit]
Description=Foo
Description=abc
AssertFileNotEmpty=1

