# tests/auto/InvalidValue/2.31/bad/Unit_RefuseManualStart_bar.target
[Unit]
Description=Foo
Description=abc
RefuseManualStart=bar

