# tests/auto/InvalidValue/2.31/bad/Unit_RefuseManualStart_foo.target
[Unit]
Description=Foo
Description=abc
RefuseManualStart=foo

