# tests/auto/InvalidValue/2.31/bad/Unit_StartLimitIntervalSec_10000000000000.target
[Unit]
Description=Foo
Description=abc
StartLimitIntervalSec=10000000000000

