# tests/auto/InvalidValue/2.31/bad/Unit_AllowIsolate_bar.target
[Unit]
Description=Foo
Description=abc
AllowIsolate=bar

