# tests/auto/InvalidValue/2.31/bad/Unit_AssertPathIsMountPoint_True.target
[Unit]
Description=Foo
Description=abc
AssertPathIsMountPoint=True

