# tests/auto/InvalidValue/2.31/bad/Unit_AssertFileIsExecutable_True.target
[Unit]
Description=Foo
Description=abc
AssertFileIsExecutable=True

