# tests/auto/InvalidValue/2.31/bad/Unit_ConditionVirtualization_bar.target
[Unit]
Description=Foo
Description=abc
ConditionVirtualization=bar

