# tests/auto/InvalidValue/2.31/bad/Unit_AssertFileIsExecutable_tmp.target
[Unit]
Description=Foo
Description=abc
AssertFileIsExecutable=tmp

