# tests/auto/InvalidValue/2.31/bad/Unit_ConditionPathIsMountPoint_1.target
[Unit]
Description=Foo
Description=abc
ConditionPathIsMountPoint=1

