# tests/auto/InvalidValue/2.31/bad/Unit_AssertDirectoryNotEmpty_1.target
[Unit]
Description=Foo
Description=abc
AssertDirectoryNotEmpty=1

