# tests/auto/InvalidValue/2.31/bad/VLAN_ReorderHeader_bar.target
[Unit]
Description=Foo
[VLAN]
ReorderHeader=bar

