# tests/auto/InvalidValue/2.31/bad/VLAN_MVRP_foo.target
[Unit]
Description=Foo
[VLAN]
MVRP=foo

