# tests/auto/InvalidValue/2.31/bad/VLAN_ReorderHeader_baz.target
[Unit]
Description=Foo
[VLAN]
ReorderHeader=baz

