# tests/auto/InvalidValue/2.31/bad/VLAN_GVRP_foo.target
[Unit]
Description=Foo
[VLAN]
GVRP=foo

