# tests/auto/InvalidValue/2.31/bad/VLAN_ReorderHeader_foo.target
[Unit]
Description=Foo
[VLAN]
ReorderHeader=foo

