# tests/auto/InvalidValue/2.31/bad/VLAN_MVRP_baz.target
[Unit]
Description=Foo
[VLAN]
MVRP=baz

