# tests/auto/InvalidValue/2.31/bad/VXLAN_Local_0--abc.target
[Unit]
Description=Foo
[VXLAN]
Local=0/abc

