# tests/auto/InvalidValue/2.31/bad/VXLAN_GroupPolicyExtension_foo.target
[Unit]
Description=Foo
[VXLAN]
GroupPolicyExtension=foo

