# tests/auto/InvalidValue/2.31/bad/VXLAN_FDBAgeingSec_-1.target
[Unit]
Description=Foo
[VXLAN]
FDBAgeingSec=-1

