# tests/auto/InvalidValue/2.31/bad/VXLAN_UDPChecksum_bar.target
[Unit]
Description=Foo
[VXLAN]
UDPChecksum=bar

