# tests/auto/InvalidValue/2.31/bad/VXLAN_DestinationPort_0.target
[Unit]
Description=Foo
[VXLAN]
DestinationPort=0

