# tests/auto/InvalidValue/2.31/bad/VXLAN_RemoteChecksumRx_bar.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumRx=bar

