# tests/auto/InvalidValue/2.31/bad/VXLAN_RemoteChecksumTx_foo.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumTx=foo

