# tests/auto/InvalidValue/2.31/bad/VXLAN_RemoteChecksumRx_foo.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumRx=foo

