# tests/auto/InvalidValue/2.31/bad/VXLAN_FlowLabel_1048576.target
[Unit]
Description=Foo
[VXLAN]
FlowLabel=1048576

