# tests/auto/InvalidValue/2.31/bad/VXLAN_RouteShortCircuit_bar.target
[Unit]
Description=Foo
[VXLAN]
RouteShortCircuit=bar

