# tests/auto/InvalidValue/2.31/bad/VXLAN_DestinationPort_65535.target
[Unit]
Description=Foo
[VXLAN]
DestinationPort=65535

