# tests/auto/InvalidValue/2.31/bad/VXLAN_FlowLabel_-1.target
[Unit]
Description=Foo
[VXLAN]
FlowLabel=-1

