# tests/auto/InvalidValue/2.31/bad/VXLAN_MaximumFDBEntries_10000000000000.target
[Unit]
Description=Foo
[VXLAN]
MaximumFDBEntries=10000000000000

