# tests/auto/InvalidValue/2.31/bad/VXLAN_RouteShortCircuit_baz.target
[Unit]
Description=Foo
[VXLAN]
RouteShortCircuit=baz

