# tests/auto/InvalidValue/2.31/bad/VXLAN_TTL_256.target
[Unit]
Description=Foo
[VXLAN]
TTL=256

