# tests/auto/InvalidValue/2.31/bad/VXLAN_RouteShortCircuit_foo.target
[Unit]
Description=Foo
[VXLAN]
RouteShortCircuit=foo

