# tests/auto/InvalidValue/2.31/bad/WireGuard_PrivateKeyFile_1.target
[Unit]
Description=Foo
[WireGuard]
PrivateKeyFile=1

