# tests/auto/InvalidValue/2.31/bad/WireGuard_PrivateKeyFile_True.target
[Unit]
Description=Foo
[WireGuard]
PrivateKeyFile=True

