# tests/auto/InvalidValue/2.31/bad/WireGuard_PrivateKeyFile_..--..--.target
[Unit]
Description=Foo
[WireGuard]
PrivateKeyFile=../../

