# tests/auto/InvalidValue/2.31/bad/WireGuard_ListenPort_0.target
[Unit]
Description=Foo
[WireGuard]
ListenPort=0

