# tests/auto/InvalidValue/2.31/bad/WireGuard_ListenPort_65535.target
[Unit]
Description=Foo
[WireGuard]
ListenPort=65535

