# tests/auto/InvalidValue/2.31/bad/WireGuardPeer_PersistentKeepalive_0.target
[Unit]
Description=Foo
[WireGuardPeer]
PersistentKeepalive=0

