# tests/auto/OptionTooNew/2.38/bad/Service_NoNewPrivileges_yes.target
[Unit]
Description=Foo
[Service]
NoNewPrivileges=123

