# tests/auto/OptionTooNew/2.38/bad/Socket_NoNewPrivileges_yes.target
[Unit]
Description=Foo
[Socket]
NoNewPrivileges=123

