# tests/auto/OptionTooNew/2.38/bad/Socket_PrivateMounts_yes.target
[Unit]
Description=Foo
[Socket]
PrivateMounts=123

