# systemdlint/tests/auto/OptionTooNew/2.48/good/Link_ReceiveQueues_1.target
[Unit]
Description=Foo
[Link]
ReceiveQueues=123

