# systemdlint/tests/auto/OptionTooNew/2.48/good/Scope_ManagedOOMPreference_none.target
[Unit]
Description=Foo
[Scope]
ManagedOOMPreference=123

