# systemdlint/tests/auto/OptionTooNew/2.48/good/Swap_NoExecPaths_abc.target
[Unit]
Description=Foo
[Swap]
NoExecPaths=123

