# systemdlint/tests/auto/OptionTooNew/2.48/good/VLAN_EgressQOSMaps_abc.target
[Unit]
Description=Foo
[VLAN]
EgressQOSMaps=123

