# tests/auto/OptionTooNew/2.16/good/Socket_KeepAliveTimeSec_0.target
[Unit]
Description=Foo
[Socket]
KeepAliveTimeSec=123

