# tests/auto/OptionTooNew/2.16/good/Socket_KeepAliveProbes_0.target
[Unit]
Description=Foo
[Socket]
KeepAliveProbes=123

