# systemdlint/tests/auto/InvalidValue/2.46/bad/CAKE_OverheadBytes_-65.target
[Unit]
Description=Foo
[CAKE]
OverheadBytes=-65

