# systemdlint/tests/auto/InvalidValue/2.46/bad/CAKE_OverheadBytes_257.target
[Unit]
Description=Foo
[CAKE]
OverheadBytes=257

