# systemdlint/tests/auto/InvalidValue/2.46/bad/CAN_FDNonISO_baz.target
[Unit]
Description=Foo
[CAN]
FDNonISO=baz

