# systemdlint/tests/auto/InvalidValue/2.46/bad/CAN_FDMode_foo.target
[Unit]
Description=Foo
[CAN]
FDMode=foo

