# systemdlint/tests/auto/InvalidValue/2.46/bad/CAN_Termination_bar.target
[Unit]
Description=Foo
[CAN]
Termination=bar

