# systemdlint/tests/auto/InvalidValue/2.46/bad/CAN_ListenOnly_bar.target
[Unit]
Description=Foo
[CAN]
ListenOnly=bar

