# systemdlint/tests/auto/InvalidValue/2.46/bad/DeficitRoundRobinScheduler_Handle_0x0.target
[Unit]
Description=Foo
[DeficitRoundRobinScheduler]
Handle=0x0

