# systemdlint/tests/auto/InvalidValue/2.46/bad/DeficitRoundRobinSchedulerClass_QuantumBytes_-1K.target
[Unit]
Description=Foo
[DeficitRoundRobinSchedulerClass]
QuantumBytes=-1K

