# systemdlint/tests/auto/InvalidValue/2.46/bad/DeficitRoundRobinSchedulerClass_QuantumBytes_10000000000000K.target
[Unit]
Description=Foo
[DeficitRoundRobinSchedulerClass]
QuantumBytes=10000000000000K

