# systemdlint/tests/auto/InvalidValue/2.46/bad/DHCPv6PrefixDelegation_Assign_foo.target
[Unit]
Description=Foo
[DHCPv6PrefixDelegation]
Assign=foo

