# systemdlint/tests/auto/InvalidValue/2.46/bad/DHCPPrefixDelegation_SubnetId_0x8000000000000000.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
SubnetId=0x8000000000000000

