# systemdlint/tests/auto/InvalidValue/2.46/bad/DHCPServer_EmitLPR_baz.target
[Unit]
Description=Foo
[DHCPServer]
EmitLPR=baz

