# systemdlint/tests/auto/InvalidValue/2.46/bad/DHCPServer_EmitSMTP_bar.target
[Unit]
Description=Foo
[DHCPServer]
EmitSMTP=bar

