# systemdlint/tests/auto/InvalidValue/2.46/bad/DHCPv6PrefixDelegation_SubnetId_0x8000000000000000.target
[Unit]
Description=Foo
[DHCPv6PrefixDelegation]
SubnetId=0x8000000000000000

