# systemdlint/tests/auto/InvalidValue/2.46/bad/EnhancedTransmissionSelection_Bands_0.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
Bands=0

