# systemdlint/tests/auto/InvalidValue/2.46/bad/EnhancedTransmissionSelection_PriorityMap_0%209999999999999%20-1.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
PriorityMap=0 9999999999999 -1

