# systemdlint/tests/auto/InvalidValue/2.46/bad/EnhancedTransmissionSelection_StrictBands_0.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
StrictBands=0

