# systemdlint/tests/auto/InvalidValue/2.46/bad/EnhancedTransmissionSelection_StrictBands_17.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
StrictBands=17

