# systemdlint/tests/auto/InvalidValue/2.46/bad/GenericRandomEarlyDetection_VirtualQueues_17.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
VirtualQueues=17

