# systemdlint/tests/auto/InvalidValue/2.46/bad/GenericRandomEarlyDetection_GenericRIO_bar.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
GenericRIO=bar

