# systemdlint/tests/auto/InvalidValue/2.46/bad/GenericRandomEarlyDetection_VirtualQueues_0.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
VirtualQueues=0

