# systemdlint/tests/auto/InvalidValue/2.46/bad/GenericRandomEarlyDetection_GenericRIO_baz.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
GenericRIO=baz

