# systemdlint/tests/auto/InvalidValue/2.46/bad/HeavyHitterFilter_Handle_0x0.target
[Unit]
Description=Foo
[HeavyHitterFilter]
Handle=0x0

