# systemdlint/tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_BufferBytes_10000000000000M.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
BufferBytes=10000000000000M

