# systemdlint/tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_CeilRate_-1M.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
CeilRate=-1M

