# systemdlint/tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_Rate_-1M.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
Rate=-1M

