# systemdlint/tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucket_RateToQuantum_-1.target
[Unit]
Description=Foo
[HierarchyTokenBucket]
RateToQuantum=-1

