# systemdlint/tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_MTUBytes_-1G.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
MTUBytes=-1G

