# systemdlint/tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_QuantumBytes_-1G.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
QuantumBytes=-1G

