# systemdlint/tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucket_RateToQuantum_10000000000000.target
[Unit]
Description=Foo
[HierarchyTokenBucket]
RateToQuantum=10000000000000

