# systemdlint/tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_MTUBytes_10000000000000M.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
MTUBytes=10000000000000M

