# systemdlint/tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_QuantumBytes_10000000000000M.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
QuantumBytes=10000000000000M

