# systemdlint/tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_BufferBytes_-1G.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
BufferBytes=-1G

