# systemdlint/tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_Rate_10000000000000G.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
Rate=10000000000000G

