# systemdlint/tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_Priority_10000000000000.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
Priority=10000000000000

