# systemdlint/tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_OverheadBytes_10000000000000M.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
OverheadBytes=10000000000000M

