# systemdlint/tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucket_DefaultClass_-0x1.target
[Unit]
Description=Foo
[HierarchyTokenBucket]
DefaultClass=-0x1

