# systemdlint/tests/auto/InvalidValue/2.46/bad/IPv6Prefix_Assign_bar.target
[Unit]
Description=Foo
[IPv6Prefix]
Assign=bar

