# systemdlint/tests/auto/InvalidValue/2.46/bad/Link_RxFlowControl_bar.target
[Unit]
Description=Foo
[Link]
RxFlowControl=bar

