# systemdlint/tests/auto/InvalidValue/2.46/bad/Link_AutoNegotiationFlowControl_foo.target
[Unit]
Description=Foo
[Link]
AutoNegotiationFlowControl=foo

