# systemdlint/tests/auto/InvalidValue/2.46/bad/Link_Group_4294967295.target
[Unit]
Description=Foo
[Link]
MACAddress=abc
Group=4294967295

