# systemdlint/tests/auto/InvalidValue/2.46/bad/Link_Group_-1.target
[Unit]
Description=Foo
[Link]
Group=-1

