# systemdlint/tests/auto/InvalidValue/2.46/bad/PFIFO_Handle_0x10000.target
[Unit]
Description=Foo
[PFIFO]
Handle=0x10000

