# systemdlint/tests/auto/InvalidValue/2.46/bad/QuickFairQueueingClass_MaxPacketBytes_-1G.target
[Unit]
Description=Foo
[QuickFairQueueingClass]
MaxPacketBytes=-1G

