# systemdlint/tests/auto/InvalidValue/2.46/bad/QuickFairQueueingClass_MaxPacketBytes_10000000000000M.target
[Unit]
Description=Foo
[QuickFairQueueingClass]
MaxPacketBytes=10000000000000M

