# systemdlint/tests/auto/InvalidValue/2.46/bad/QuickFairQueueingClass_Weight_0.target
[Unit]
Description=Foo
[QuickFairQueueingClass]
Weight=0

