# systemdlint/tests/auto/InvalidValue/2.46/bad/Socket_PassPacketInfo_bar.target
[Unit]
Description=Foo
[Socket]
PassPacketInfo=bar

