# systemdlint/tests/auto/InvalidValue/2.46/bad/SR-IOV_QualityOfService_0.target
[Unit]
Description=Foo
[SR-IOV]
QualityOfService=0

