# systemdlint/tests/auto/InvalidValue/2.46/bad/StochasticFairBlue_PacketLimit_4294967295.target
[Unit]
Description=Foo
[StochasticFairBlue]
PacketLimit=4294967295

