# systemdlint/tests/auto/InvalidValue/2.46/bad/Unit_AssertEnvironment_a%3Eb.target
[Unit]
Description=Foo
Description=abc
AssertEnvironment=a>b

