# systemdlint/tests/auto/InvalidValue/2.43/good/BridgeFDB_VNI_16777215.target
[Unit]
Description=Foo
[BridgeFDB]
MACAddress=01:23:45:67:89:ab
VNI=16777215

