# systemdlint/tests/auto/InvalidValue/2.43/good/Bridge_ProxyARP_false.target
[Unit]
Description=Foo
[Bridge]
ProxyARP=false

