# systemdlint/tests/auto/InvalidValue/2.43/good/Bridge_ProxyARP_no.target
[Unit]
Description=Foo
[Bridge]
ProxyARP=no

