# systemdlint/tests/auto/InvalidValue/2.43/good/Bridge_ProxyARPWiFi_on.target
[Unit]
Description=Foo
[Bridge]
ProxyARPWiFi=on

