# systemdlint/tests/auto/InvalidValue/2.43/good/BridgeFDB_AssociatedWith_router.target
[Unit]
Description=Foo
[BridgeFDB]
MACAddress=abc
AssociatedWith=router

