# systemdlint/tests/auto/InvalidValue/2.43/good/Bridge_ProxyARP_on.target
[Unit]
Description=Foo
[Bridge]
ProxyARP=on

