# systemdlint/tests/auto/InvalidValue/2.43/good/Bridge_ProxyARP_true.target
[Unit]
Description=Foo
[Bridge]
ProxyARP=true

