# systemdlint/tests/auto/InvalidValue/2.43/good/Bridge_ProxyARPWiFi_1.target
[Unit]
Description=Foo
[Bridge]
ProxyARPWiFi=1

