# systemdlint/tests/auto/InvalidValue/2.43/good/Bridge_ProxyARP_yes.target
[Unit]
Description=Foo
[Bridge]
ProxyARP=yes

