# systemdlint/tests/auto/InvalidValue/2.43/good/Bridge_ProxyARPWiFi_false.target
[Unit]
Description=Foo
[Bridge]
ProxyARPWiFi=false

