# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv6_ListenPort_1.target
[Unit]
Description=Foo
[DHCPv6]
ListenPort=1

