# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv6_UseTimezone_off.target
[Unit]
Description=Foo
[DHCPv6]
UseTimezone=off

