# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv6_RouteMetric_9999999999999.target
[Unit]
Description=Foo
[DHCPv6]
RouteMetric=9999999999999

