# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv4_ClientIdentifier_duid-only.target
[Unit]
Description=Foo
[DHCPv4]
ClientIdentifier=duid-only

