# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv6_Anonymize_on.target
[Unit]
Description=Foo
[DHCPv6]
Anonymize=on

