# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv6_UseHostname_1.target
[Unit]
Description=Foo
[DHCPv6]
UseHostname=1

