# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv4_UseHostname_false.target
[Unit]
Description=Foo
[DHCPv4]
UseHostname=false

