# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv4_UseMTU_1.target
[Unit]
Description=Foo
[DHCPv4]
UseMTU=1

