# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv6_ForceDHCPv6PDOtherInformation_1.target
[Unit]
Description=Foo
[DHCPv6]
ForceDHCPv6PDOtherInformation=1

