# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv6_UseDNS_1.target
[Unit]
Description=Foo
[DHCPv6]
UseDNS=1

