# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv6_UseDNS_on.target
[Unit]
Description=Foo
[DHCPv6]
UseDNS=on

