# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv4_IAID_0.target
[Unit]
Description=Foo
[DHCPv4]
IAID=0

