# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv6_UseHostname_no.target
[Unit]
Description=Foo
[DHCPv6]
UseHostname=no

