# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv6_Hostname_abc.target
[Unit]
Description=Foo
[DHCPv6]
Hostname=abc

