# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv6_RouteTable_4294967295.target
[Unit]
Description=Foo
[DHCPv6]
RouteTable=4294967295

