# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv4_UseDNS_off.target
[Unit]
Description=Foo
[DHCPv4]
UseDNS=off

