# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv4_RapidCommit_0.target
[Unit]
Description=Foo
[DHCPv4]
RapidCommit=0

