# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv6_KeepConfiguration_abc.target
[Unit]
Description=Foo
[DHCPv6]
KeepConfiguration=abc

