# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv4_UseDomains_off.target
[Unit]
Description=Foo
[DHCPv4]
UseDomains=off

