# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv4_UseTimezone_true.target
[Unit]
Description=Foo
[DHCPv4]
UseTimezone=true

