# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv6_UseDomains_1.target
[Unit]
Description=Foo
[DHCPv6]
UseDomains=1

