# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv4_ListenPort_65535.target
[Unit]
Description=Foo
[DHCPv4]
ListenPort=65535

