# systemdlint/tests/auto/InvalidValue/2.43/good/DHCP_KeepConfiguration_abc.target
[Unit]
Description=Foo
[DHCP]
KeepConfiguration=abc

