# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv6_UseRoutes_off.target
[Unit]
Description=Foo
[DHCPv6]
UseRoutes=off

