# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv6_IAID_9999999999999.target
[Unit]
Description=Foo
[DHCPv6]
IAID=9999999999999

