# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv6_DUIDRawData_abc.target
[Unit]
Description=Foo
[DHCPv6]
DUIDRawData=abc

