# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv6_UseNTP_no.target
[Unit]
Description=Foo
[DHCPv6]
UseNTP=no

