# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv6_ClientIdentifier_duid-only.target
[Unit]
Description=Foo
[DHCPv6]
ClientIdentifier=duid-only

