# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv4_UseTimezone_0.target
[Unit]
Description=Foo
[DHCPv4]
UseTimezone=0

