# systemdlint/tests/auto/InvalidValue/2.43/good/DHCPv4_Hostname_abc.target
[Unit]
Description=Foo
[DHCPv4]
Hostname=abc

