# systemdlint/tests/auto/InvalidValue/2.43/good/FooOverUDP_Peer_1.1.1.1.target
[Unit]
Description=Foo
[FooOverUDP]
Peer=1.1.1.1

