# systemdlint/tests/auto/InvalidValue/2.43/good/GENEVE_GenericProtocolExtension_on.target
[Unit]
Description=Foo
[GENEVE]
GenericProtocolExtension=on

