# systemdlint/tests/auto/InvalidValue/2.43/good/GENEVE_IPDoNotFragment_false.target
[Unit]
Description=Foo
[GENEVE]
IPDoNotFragment=false

