# systemdlint/tests/auto/InvalidValue/2.43/good/GENEVE_GenericProtocolExtension_1.target
[Unit]
Description=Foo
[GENEVE]
GenericProtocolExtension=1

