# systemdlint/tests/auto/InvalidValue/2.43/good/GENEVE_IPDoNotFragment_on.target
[Unit]
Description=Foo
[GENEVE]
IPDoNotFragment=on

