# systemdlint/tests/auto/InvalidValue/2.43/good/GENEVE_IPDoNotFragment_no.target
[Unit]
Description=Foo
[GENEVE]
IPDoNotFragment=no

