# systemdlint/tests/auto/InvalidValue/2.43/good/GENEVE_IPDoNotFragment_off.target
[Unit]
Description=Foo
[GENEVE]
IPDoNotFragment=off

