# systemdlint/tests/auto/InvalidValue/2.43/good/Link_OriginalName_abc.target
[Unit]
Description=Foo
[Link]
OriginalName=abc

