# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_4ms.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=4ms

