# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_11y.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=11y

