# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_8weeks.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=8weeks

