# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_21usec.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=21usec

