# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_26months.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=26months

