# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_20usec.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=20usec

