# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_9min.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=9min

