# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_0sec.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=0sec

