# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_17second.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=17second

