# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_21usec.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=21usec

