# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_20w.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=20w

