# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_28%C2%B5s.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=28µs

