# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_9month.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=9month

