# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_7seconds.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=7seconds

