# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_18sec.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=18sec

