# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_27msec.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=27msec

