# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_12second.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=12second

