# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_15min.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=15min

