# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_14M.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=14M

