# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_24years.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=24years

