# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_17minutes.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=17minutes

