# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_25s.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=25s

