# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_14msec.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=14msec

