# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_11y.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=11y

