# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_24w.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=24w

