# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_NUMAPolicy_default.target
[Unit]
Description=Foo
[Manager]
NUMAPolicy=default

