# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_22hr.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=22hr

