# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_2us.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=2us

