# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_9month.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=9month

