# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_24w.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=24w

