# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_23m.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=23m

