# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_15ms.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=15ms

