# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_16week.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=16week

