# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_27y.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=27y

