# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_NUMAPolicy_preferred.target
[Unit]
Description=Foo
[Manager]
NUMAPolicy=preferred

