# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_5hr.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=5hr

