# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_4ms.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=4ms

