# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_19d.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=19d

