# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_17minutes.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=17minutes

