# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_23m.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=23m

