# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_16week.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=16week

