# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_5hour.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=5hour

