# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_3hours.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=3hours

