# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_15min.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=15min

