# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_2us.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=2us

