# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_KExecWatchdogSec_19d.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=19d

