# systemdlint/tests/auto/InvalidValue/2.43/good/Manager_RebootWatchdogSec_12month.target
[Unit]
Description=Foo
[Manager]
RebootWatchdogSec=12month

