# systemdlint/tests/auto/InvalidValue/2.43/good/Mount_DefaultMemoryLow_9999999999999B.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
DefaultMemoryLow=9999999999999B

