# systemdlint/tests/auto/InvalidValue/2.43/good/Neighbor_LinkLayerAddress_abc.target
[Unit]
Description=Foo
[Neighbor]
LinkLayerAddress=abc

