# systemdlint/tests/auto/InvalidValue/2.43/good/Network_DefaultRouteOnDevice_off.target
[Unit]
Description=Foo
[Network]
DefaultRouteOnDevice=off

