# systemdlint/tests/auto/InvalidValue/2.43/good/Network_DefaultRouteOnDevice_1.target
[Unit]
Description=Foo
[Network]
DefaultRouteOnDevice=1

