# systemdlint/tests/auto/InvalidValue/2.43/good/Network_DefaultRouteOnDevice_no.target
[Unit]
Description=Foo
[Network]
DefaultRouteOnDevice=no

