# systemdlint/tests/auto/InvalidValue/2.43/good/Network_DefaultRouteOnDevice_true.target
[Unit]
Description=Foo
[Network]
DefaultRouteOnDevice=true

