# systemdlint/tests/auto/InvalidValue/2.43/good/Route_FastOpenNoCookie_yes.target
[Unit]
Description=Foo
[Route]
FastOpenNoCookie=yes

