# systemdlint/tests/auto/InvalidValue/2.43/good/Route_TTLPropagate_off.target
[Unit]
Description=Foo
[Route]
TTLPropagate=off

