# systemdlint/tests/auto/InvalidValue/2.43/good/Route_FastOpenNoCookie_true.target
[Unit]
Description=Foo
[Route]
FastOpenNoCookie=true

