# systemdlint/tests/auto/InvalidValue/2.43/good/Route_TTLPropagate_no.target
[Unit]
Description=Foo
[Route]
TTLPropagate=no

