# systemdlint/tests/auto/InvalidValue/2.43/good/Service_TimeoutAbortSec_5hour.target
[Unit]
Description=Foo
[Service]
TimeoutAbortSec=5hour

