# systemdlint/tests/auto/InvalidValue/2.43/good/Service_TimeoutAbortSec_10day.target
[Unit]
Description=Foo
[Service]
TimeoutAbortSec=10day

