# systemdlint/tests/auto/InvalidValue/2.43/good/Service_TimeoutAbortSec_4ms.target
[Unit]
Description=Foo
[Service]
TimeoutAbortSec=4ms

