# systemdlint/tests/auto/InvalidValue/2.43/good/Service_NUMAPolicy_preferred.target
[Unit]
Description=Foo
[Service]
NUMAPolicy=preferred

