# systemdlint/tests/auto/InvalidValue/2.43/good/Service_NUMAPolicy_bind.target
[Unit]
Description=Foo
[Service]
NUMAPolicy=bind

