# systemdlint/tests/auto/InvalidValue/2.43/good/Service_TimeoutAbortSec_6days.target
[Unit]
Description=Foo
[Service]
TimeoutAbortSec=6days

