# systemdlint/tests/auto/InvalidValue/2.43/good/Service_TimeoutAbortSec_8weeks.target
[Unit]
Description=Foo
[Service]
TimeoutAbortSec=8weeks

