# systemdlint/tests/auto/InvalidValue/2.43/good/Service_TimeoutAbortSec_23m.target
[Unit]
Description=Foo
[Service]
TimeoutAbortSec=23m

