# systemdlint/tests/auto/InvalidValue/2.43/good/Service_TimeoutAbortSec_13minute.target
[Unit]
Description=Foo
[Service]
TimeoutAbortSec=13minute

