# systemdlint/tests/auto/InvalidValue/2.43/good/Service_TimeoutAbortSec_26months.target
[Unit]
Description=Foo
[Service]
TimeoutAbortSec=26months

