# systemdlint/tests/auto/InvalidValue/2.43/good/Service_TimeoutAbortSec_7minutes.target
[Unit]
Description=Foo
[Service]
TimeoutAbortSec=7minutes

