# systemdlint/tests/auto/InvalidValue/2.43/good/Service_TimeoutAbortSec_6h.target
[Unit]
Description=Foo
[Service]
TimeoutAbortSec=6h

