# systemdlint/tests/auto/InvalidValue/2.43/good/Service_DefaultMemoryMin_9999999999999T.target
[Unit]
Description=Foo
[Service]
DefaultMemoryMin=9999999999999T

