# systemdlint/tests/auto/InvalidValue/2.43/good/Service_NUMAPolicy_local.target
[Unit]
Description=Foo
[Service]
NUMAPolicy=local

