# systemdlint/tests/auto/InvalidValue/2.43/good/Service_TimeoutAbortSec_1day.target
[Unit]
Description=Foo
[Service]
TimeoutAbortSec=1day

