# systemdlint/tests/auto/InvalidValue/2.43/good/Service_TimeoutAbortSec_2d.target
[Unit]
Description=Foo
[Service]
TimeoutAbortSec=2d

