# systemdlint/tests/auto/InvalidValue/2.43/good/Service_ExecCondition_abc.target
[Unit]
Description=Foo
[Service]
ExecCondition=abc

