# systemdlint/tests/auto/InvalidValue/2.43/good/Service_TimeoutAbortSec_27msec.target
[Unit]
Description=Foo
[Service]
TimeoutAbortSec=27msec

