# systemdlint/tests/auto/InvalidValue/2.43/good/Service_TimeoutAbortSec_20usec.target
[Unit]
Description=Foo
[Service]
TimeoutAbortSec=20usec

