# systemdlint/tests/auto/InvalidValue/2.43/good/Service_DefaultMemoryMin_0G.target
[Unit]
Description=Foo
[Service]
DefaultMemoryMin=0G

