# systemdlint/tests/auto/InvalidValue/2.43/good/Socket_DefaultMemoryLow_9999999999999B.target
[Unit]
Description=Foo
[Socket]
DefaultMemoryLow=9999999999999B

