# systemdlint/tests/auto/InvalidValue/2.43/good/Socket_NUMAPolicy_interleave.target
[Unit]
Description=Foo
[Socket]
NUMAPolicy=interleave

