# systemdlint/tests/auto/InvalidValue/2.43/good/Socket_NUMAPolicy_local.target
[Unit]
Description=Foo
[Socket]
NUMAPolicy=local

