# systemdlint/tests/auto/InvalidValue/2.43/good/Swap_NUMAPolicy_local.target
[Unit]
Description=Foo
[Swap]
NUMAPolicy=local

