# systemdlint/tests/auto/InvalidValue/2.43/good/Swap_NUMAPolicy_bind.target
[Unit]
Description=Foo
[Swap]
NUMAPolicy=bind

