# systemdlint/tests/auto/InvalidValue/2.43/good/Tunnel_AssignToLoopback_false.target
[Unit]
Description=Foo
[Tunnel]
AssignToLoopback=false

