# systemdlint/tests/auto/InvalidValue/2.43/good/Tunnel_AssignToLoopback_true.target
[Unit]
Description=Foo
[Tunnel]
AssignToLoopback=true

