# systemdlint/tests/auto/InvalidValue/2.43/good/VXLAN_IPDoNotFragment_false.target
[Unit]
Description=Foo
[VXLAN]
IPDoNotFragment=false

