# systemdlint/tests/auto/InvalidValue/2.43/good/VXLAN_IPDoNotFragment_yes.target
[Unit]
Description=Foo
[VXLAN]
IPDoNotFragment=yes

