# systemdlint/tests/auto/InvalidValue/2.43/good/VXLAN_GenericProtocolExtension_true.target
[Unit]
Description=Foo
[VXLAN]
GenericProtocolExtension=true

