# systemdlint/tests/auto/InvalidValue/2.43/good/VXLAN_IPDoNotFragment_1.target
[Unit]
Description=Foo
[VXLAN]
IPDoNotFragment=1

