# systemdlint/tests/auto/InvalidValue/2.43/good/VXLAN_GenericProtocolExtension_false.target
[Unit]
Description=Foo
[VXLAN]
GenericProtocolExtension=false

