# systemdlint/tests/auto/InvalidValue/2.43/good/VXLAN_IPDoNotFragment_true.target
[Unit]
Description=Foo
[VXLAN]
IPDoNotFragment=true

