# systemdlint/tests/auto/InvalidValue/2.43/good/VXLAN_GenericProtocolExtension_on.target
[Unit]
Description=Foo
[VXLAN]
GenericProtocolExtension=on

