# tests/auto/OptionTooNew/2.32/bad/Mount_RestrictNamespaces_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
RestrictNamespaces=123

