# tests/auto/OptionTooNew/2.32/bad/Socket_BindReadOnlyPaths_abc.target
[Unit]
Description=Foo
[Socket]
BindReadOnlyPaths=123

