# tests/auto/OptionTooNew/2.32/bad/Socket_RestrictNamespaces_abc.target
[Unit]
Description=Foo
[Socket]
RestrictNamespaces=123

