# systemdlint/tests/auto/InvalidValue/2.12/good/Service_ControlGroup_abc.target
[Unit]
Description=Foo
[Service]
ControlGroup=abc

