# systemdlint/tests/auto/InvalidValue/2.12/good/Socket_ControlGroupAttribute_abc.target
[Unit]
Description=Foo
[Socket]
ControlGroupAttribute=abc

