# systemdlint/tests/auto/InvalidValue/2.12/good/Swap_ControlGroupPersistent_abc.target
[Unit]
Description=Foo
[Swap]
ControlGroupPersistent=abc

