# systemdlint/tests/auto/InvalidValue/2.12/good/Swap_ControlGroupAttribute_abc.target
[Unit]
Description=Foo
[Swap]
ControlGroupAttribute=abc

