# tests/auto/OptionTooNew/2.45/bad/CAKE_Bandwidth_0K.target
[Unit]
Description=Foo
[CAKE]
Bandwidth=123

