# tests/auto/OptionTooNew/2.45/bad/CAKE_OverheadBytes_-64.target
[Unit]
Description=Foo
[CAKE]
OverheadBytes=123

