# tests/auto/OptionTooNew/2.45/bad/CAN_DataSamplePoint_abc.target
[Unit]
Description=Foo
[CAN]
DataSamplePoint=123

