# tests/auto/OptionTooNew/2.45/bad/CAN_BitRate_1K.target
[Unit]
Description=Foo
[CAN]
BitRate=123

