# tests/auto/OptionTooNew/2.45/bad/CAN_FDNonISO_yes.target
[Unit]
Description=Foo
[CAN]
FDNonISO=123

