# tests/auto/OptionTooNew/2.45/bad/DeficitRoundRobinSchedulerClass_ClassId_abc.target
[Unit]
Description=Foo
[DeficitRoundRobinSchedulerClass]
ClassId=123

