# tests/auto/OptionTooNew/2.45/bad/DeficitRoundRobinScheduler_Handle_0x1.target
[Unit]
Description=Foo
[DeficitRoundRobinScheduler]
Handle=123

