# tests/auto/OptionTooNew/2.45/bad/DHCPv6_SendOption_abc.target
[Unit]
Description=Foo
[DHCPv6]
SendOption=123

