# tests/auto/OptionTooNew/2.45/bad/DHCPServer_EmitSMTP_yes.target
[Unit]
Description=Foo
[DHCPServer]
EmitSMTP=123

