# tests/auto/OptionTooNew/2.45/bad/DHCPv4_DenyList_abc.target
[Unit]
Description=Foo
[DHCPv4]
DenyList=123

