# tests/auto/OptionTooNew/2.45/bad/DHCPv6PrefixDelegation_Token_abc.target
[Unit]
Description=Foo
[DHCPv6PrefixDelegation]
Token=123

