# tests/auto/OptionTooNew/2.45/bad/EnhancedTransmissionSelection_StrictBands_1.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
StrictBands=123

