# tests/auto/OptionTooNew/2.45/bad/EnhancedTransmissionSelection_Bands_1.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
Bands=123

