# tests/auto/OptionTooNew/2.45/bad/HierarchyTokenBucketClass_Parent_abc.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
Parent=123

