# tests/auto/OptionTooNew/2.45/bad/Link_TxFlowControl_yes.target
[Unit]
Description=Foo
[Link]
TxFlowControl=123

