# tests/auto/OptionTooNew/2.45/bad/Network_IPv6LinkLocalAddressGenerationMode_eui64.target
[Unit]
Description=Foo
[Network]
IPv6LinkLocalAddressGenerationMode=123

