# tests/auto/OptionTooNew/2.45/bad/Network_IPv4AcceptLocal_yes.target
[Unit]
Description=Foo
[Network]
IPv4AcceptLocal=123

