# tests/auto/OptionTooNew/2.45/bad/PFIFO_PacketLimit_0.target
[Unit]
Description=Foo
[PFIFO]
PacketLimit=123

