# tests/auto/OptionTooNew/2.45/bad/PFIFOHeadDrop_PacketLimit_0.target
[Unit]
Description=Foo
[PFIFOHeadDrop]
PacketLimit=123

