# tests/auto/OptionTooNew/2.45/bad/PFIFOFast_Parent_abc.target
[Unit]
Description=Foo
[PFIFOFast]
Parent=123

