# tests/auto/OptionTooNew/2.45/bad/PFIFOHeadDrop_Handle_0x1.target
[Unit]
Description=Foo
[PFIFOHeadDrop]
Handle=123

