# tests/auto/OptionTooNew/2.45/bad/QuickFairQueueing_Parent_abc.target
[Unit]
Description=Foo
[QuickFairQueueing]
Parent=123

