# tests/auto/OptionTooNew/2.45/bad/QuickFairQueueing_Handle_0x1.target
[Unit]
Description=Foo
[QuickFairQueueing]
Handle=123

