# tests/auto/OptionTooNew/2.45/bad/Service_TimeoutStopFailureMode_terminate.target
[Unit]
Description=Foo
[Service]
TimeoutStopFailureMode=123

