# tests/auto/OptionTooNew/2.45/bad/Service_RootHash_abc.target
[Unit]
Description=Foo
[Service]
RootHash=123

