# tests/auto/OptionTooNew/2.45/bad/Service_TimeoutStartFailureMode_terminate.target
[Unit]
Description=Foo
[Service]
TimeoutStartFailureMode=123

