# tests/auto/OptionTooNew/2.45/bad/Socket_CoredumpFilter_0.target
[Unit]
Description=Foo
[Socket]
CoredumpFilter=123

