# tests/auto/OptionTooNew/2.45/bad/Socket_RootHashSignature_abc.target
[Unit]
Description=Foo
[Socket]
RootHashSignature=123

