# tests/auto/OptionTooNew/2.45/bad/Swap_RootHash_abc.target
[Unit]
Description=Foo
[Swap]
RootHash=123

