# tests/auto/OptionTooNew/2.45/bad/Swap_CoredumpFilter_0.target
[Unit]
Description=Foo
[Swap]
CoredumpFilter=123

