# tests/auto/OptionDeprecated/2.46/good/IPv6PrefixDelegation_DNS_abc.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
DNS=123

