# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_FlowIsolationMode_none.target
[Unit]
Description=Foo
[CAKE]
FlowIsolationMode=none

