# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_FlowIsolationMode_flows.target
[Unit]
Description=Foo
[CAKE]
FlowIsolationMode=flows

