# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_UseRawPacketSize_off.target
[Unit]
Description=Foo
[CAKE]
UseRawPacketSize=off

