# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_NAT_true.target
[Unit]
Description=Foo
[CAKE]
NAT=true

