# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_FlowIsolationMode_hosts.target
[Unit]
Description=Foo
[CAKE]
FlowIsolationMode=hosts

