# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_SplitGSO_false.target
[Unit]
Description=Foo
[CAKE]
SplitGSO=false

