# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_FlowIsolationMode_dst-host.target
[Unit]
Description=Foo
[CAKE]
FlowIsolationMode=dst-host

