# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_CompensationMode_none.target
[Unit]
Description=Foo
[CAKE]
CompensationMode=none

