# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_AutoRateIngress_false.target
[Unit]
Description=Foo
[CAKE]
AutoRateIngress=false

