# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_SplitGSO_1.target
[Unit]
Description=Foo
[CAKE]
SplitGSO=1

