# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_CompensationMode_atm.target
[Unit]
Description=Foo
[CAKE]
CompensationMode=atm

