# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_FlowIsolationMode_src-host.target
[Unit]
Description=Foo
[CAKE]
FlowIsolationMode=src-host

