# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_FirewallMark_1.target
[Unit]
Description=Foo
[CAKE]
FirewallMark=1

