# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_PriorityQueueingPreset_besteffort.target
[Unit]
Description=Foo
[CAKE]
PriorityQueueingPreset=besteffort

