# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_UseRawPacketSize_0.target
[Unit]
Description=Foo
[CAKE]
UseRawPacketSize=0

