# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_UseRawPacketSize_yes.target
[Unit]
Description=Foo
[CAKE]
UseRawPacketSize=yes

