# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_FlowIsolationMode_dual-dst-host.target
[Unit]
Description=Foo
[CAKE]
FlowIsolationMode=dual-dst-host

