# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_Wash_off.target
[Unit]
Description=Foo
[CAKE]
Wash=off

