# systemdlint/tests/auto/InvalidValue/2.50/good/CAKE_AutoRateIngress_yes.target
[Unit]
Description=Foo
[CAKE]
AutoRateIngress=yes

