# systemdlint/tests/auto/InvalidValue/2.50/good/CAN_Loopback_false.target
[Unit]
Description=Foo
[CAN]
Loopback=false

