# systemdlint/tests/auto/InvalidValue/2.50/good/CAN_PhaseBufferSegment1_4294967295.target
[Unit]
Description=Foo
[CAN]
PhaseBufferSegment1=4294967295

