# systemdlint/tests/auto/InvalidValue/2.50/good/CAN_PhaseBufferSegment1_0.target
[Unit]
Description=Foo
[CAN]
PhaseBufferSegment1=0

