# systemdlint/tests/auto/InvalidValue/2.50/good/CAN_SyncJumpWidth_0.target
[Unit]
Description=Foo
[CAN]
SyncJumpWidth=0

