# systemdlint/tests/auto/InvalidValue/2.50/good/CAN_TimeQuantaNSec_0.target
[Unit]
Description=Foo
[CAN]
TimeQuantaNSec=0

