# systemdlint/tests/auto/InvalidValue/2.50/good/CAN_DataPhaseBufferSegment1_4294967295.target
[Unit]
Description=Foo
[CAN]
DataPhaseBufferSegment1=4294967295

