# systemdlint/tests/auto/InvalidValue/2.50/good/CAN_DataSyncJumpWidth_4294967295.target
[Unit]
Description=Foo
[CAN]
DataSyncJumpWidth=4294967295

