# systemdlint/tests/auto/InvalidValue/2.50/good/CAN_PresumeAck_abc.target
[Unit]
Description=Foo
[CAN]
PresumeAck=abc

