# systemdlint/tests/auto/InvalidValue/2.50/good/CAN_DataTimeQuantaNSec_4294967295.target
[Unit]
Description=Foo
[CAN]
DataTimeQuantaNSec=4294967295

