# systemdlint/tests/auto/InvalidValue/2.50/good/CAN_DataTimeQuantaNSec_0.target
[Unit]
Description=Foo
[CAN]
DataTimeQuantaNSec=0

