# systemdlint/tests/auto/InvalidValue/2.50/good/CAN_ClassicDataLengthCode_abc.target
[Unit]
Description=Foo
[CAN]
ClassicDataLengthCode=abc

