# systemdlint/tests/auto/InvalidValue/2.50/good/CAN_OneShot_off.target
[Unit]
Description=Foo
[CAN]
OneShot=off

