# systemdlint/tests/auto/InvalidValue/2.50/good/DHCPPrefixDelegation_Assign_false.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Assign=false

