# systemdlint/tests/auto/InvalidValue/2.50/good/DHCPPrefixDelegation_Announce_yes.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Announce=yes

