# systemdlint/tests/auto/InvalidValue/2.50/good/DHCPPrefixDelegation_Announce_false.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Announce=false

