# systemdlint/tests/auto/InvalidValue/2.50/good/DHCPv4_Label_abc.target
[Unit]
Description=Foo
[DHCPv4]
Label=abc

