# systemdlint/tests/auto/InvalidValue/2.50/good/DHCPPrefixDelegation_Assign_yes.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Assign=yes

