# systemdlint/tests/auto/InvalidValue/2.50/good/DHCPPrefixDelegation_UplinkInterface_abc.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
UplinkInterface=abc

