# systemdlint/tests/auto/InvalidValue/2.50/good/DHCPPrefixDelegation_RouteMetric_4294967295.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
RouteMetric=4294967295

