# systemdlint/tests/auto/InvalidValue/2.50/good/DHCPPrefixDelegation_Assign_true.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Assign=true

