# systemdlint/tests/auto/InvalidValue/2.50/good/DHCPPrefixDelegation_Assign_off.target
[Unit]
Description=Foo
[DHCPPrefixDelegation]
Assign=off

