# systemdlint/tests/auto/InvalidValue/2.50/good/IPoIB_IgnoreUserspaceMulticastGroup_off.target
[Unit]
Description=Foo
[IPoIB]
IgnoreUserspaceMulticastGroup=off

