# systemdlint/tests/auto/InvalidValue/2.50/good/IPoIB_IgnoreUserspaceMulticastGroup_false.target
[Unit]
Description=Foo
[IPoIB]
IgnoreUserspaceMulticastGroup=false

