# systemdlint/tests/auto/InvalidValue/2.50/good/IPv6AcceptRA_UseMTU_yes.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseMTU=yes

