# systemdlint/tests/auto/InvalidValue/2.50/good/IPv6SendRA_UplinkInterface_abc.target
[Unit]
Description=Foo
[IPv6SendRA]
UplinkInterface=abc

