# systemdlint/tests/auto/InvalidValue/2.50/good/IPv6Prefix_Token_abc.target
[Unit]
Description=Foo
[IPv6Prefix]
Token=abc

