# systemdlint/tests/auto/InvalidValue/2.50/good/IPv6AcceptRA_UseMTU_0.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseMTU=0

