# systemdlint/tests/auto/InvalidValue/2.50/good/IPv6AcceptRA_UseRoutePrefix_on.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseRoutePrefix=on

