# systemdlint/tests/auto/InvalidValue/2.50/good/Link_WakeOnLanPassword_abc.target
[Unit]
Description=Foo
[Link]
WakeOnLanPassword=abc

