# systemdlint/tests/auto/InvalidValue/2.50/good/Login_HandlePowerKeyLongPress_ignore.target
[Unit]
Description=Foo
[Login]
HandlePowerKeyLongPress=ignore

