# systemdlint/tests/auto/InvalidValue/2.50/good/Login_HandleRebootKeyLongPress_lock.target
[Unit]
Description=Foo
[Login]
HandleRebootKeyLongPress=lock

