# systemdlint/tests/auto/InvalidValue/2.50/good/Login_HandleSuspendKeyLongPress_hibernate.target
[Unit]
Description=Foo
[Login]
HandleSuspendKeyLongPress=hibernate

