# systemdlint/tests/auto/InvalidValue/2.50/good/Login_HandleRebootKeyLongPress_hibernate.target
[Unit]
Description=Foo
[Login]
HandleRebootKeyLongPress=hibernate

