# systemdlint/tests/auto/InvalidValue/2.50/good/Login_HandlePowerKeyLongPress_lock.target
[Unit]
Description=Foo
[Login]
HandlePowerKeyLongPress=lock

