# systemdlint/tests/auto/InvalidValue/2.50/good/Login_HandleSuspendKeyLongPress_lock.target
[Unit]
Description=Foo
[Login]
HandleSuspendKeyLongPress=lock

