# systemdlint/tests/auto/InvalidValue/2.50/good/Login_HandleHibernateKeyLongPress_hybrid-sleep.target
[Unit]
Description=Foo
[Login]
HandleHibernateKeyLongPress=hybrid-sleep

