# systemdlint/tests/auto/InvalidValue/2.50/good/Login_HandleRebootKeyLongPress_ignore.target
[Unit]
Description=Foo
[Login]
HandleRebootKeyLongPress=ignore

