# systemdlint/tests/auto/InvalidValue/2.50/good/Login_HandleSuspendKeyLongPress_ignore.target
[Unit]
Description=Foo
[Login]
HandleSuspendKeyLongPress=ignore

