# systemdlint/tests/auto/InvalidValue/2.50/good/Login_HandleHibernateKeyLongPress_halt.target
[Unit]
Description=Foo
[Login]
HandleHibernateKeyLongPress=halt

