# systemdlint/tests/auto/InvalidValue/2.50/good/Login_HandleRebootKeyLongPress_hybrid-sleep.target
[Unit]
Description=Foo
[Login]
HandleRebootKeyLongPress=hybrid-sleep

