# systemdlint/tests/auto/InvalidValue/2.50/good/Login_HandleHibernateKeyLongPress_suspend-then-hibernate.target
[Unit]
Description=Foo
[Login]
HandleHibernateKeyLongPress=suspend-then-hibernate

