# systemdlint/tests/auto/InvalidValue/2.50/good/Login_HandleHibernateKeyLongPress_reboot.target
[Unit]
Description=Foo
[Login]
HandleHibernateKeyLongPress=reboot

