# systemdlint/tests/auto/InvalidValue/2.50/good/Login_HandleSuspendKeyLongPress_reboot.target
[Unit]
Description=Foo
[Login]
HandleSuspendKeyLongPress=reboot

