# systemdlint/tests/auto/InvalidValue/2.50/good/Login_HandleSuspendKeyLongPress_halt.target
[Unit]
Description=Foo
[Login]
HandleSuspendKeyLongPress=halt

