# systemdlint/tests/auto/InvalidValue/2.50/good/Login_HandleSuspendKeyLongPress_hybrid-sleep.target
[Unit]
Description=Foo
[Login]
HandleSuspendKeyLongPress=hybrid-sleep

