# systemdlint/tests/auto/InvalidValue/2.50/good/Path_TriggerLimitIntervalSec_25years.target
[Unit]
Description=Foo
[Path]
TriggerLimitIntervalSec=25years

