# systemdlint/tests/auto/InvalidValue/2.50/good/Path_TriggerLimitIntervalSec_15ms.target
[Unit]
Description=Foo
[Path]
TriggerLimitIntervalSec=15ms

