# systemdlint/tests/auto/InvalidValue/2.50/good/Path_TriggerLimitIntervalSec_8minute.target
[Unit]
Description=Foo
[Path]
TriggerLimitIntervalSec=8minute

