# systemdlint/tests/auto/InvalidValue/2.50/good/Path_TriggerLimitIntervalSec_18sec.target
[Unit]
Description=Foo
[Path]
TriggerLimitIntervalSec=18sec

