# systemdlint/tests/auto/InvalidValue/2.50/good/Path_TriggerLimitIntervalSec_14msec.target
[Unit]
Description=Foo
[Path]
TriggerLimitIntervalSec=14msec

