# systemdlint/tests/auto/InvalidValue/2.50/good/Path_TriggerLimitIntervalSec_17second.target
[Unit]
Description=Foo
[Path]
TriggerLimitIntervalSec=17second

