# systemdlint/tests/auto/InvalidValue/2.50/good/Path_TriggerLimitBurst_0.target
[Unit]
Description=Foo
[Path]
TriggerLimitBurst=0

