# systemdlint/tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_25years.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=25years

