# systemdlint/tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_20usec.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=20usec

