# systemdlint/tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_4hour.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=4hour

