# systemdlint/tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_18sec.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=18sec

