# systemdlint/tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_3hours.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=3hours

