# systemdlint/tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_infinity.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=infinity

