# systemdlint/tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_10M.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=10M

