# systemdlint/tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_2d.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=2d

