# systemdlint/tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_5hr.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=5hr

