# systemdlint/tests/auto/InvalidValue/2.50/good/Scope_RuntimeRandomizedExtraSec_24w.target
[Unit]
Description=Foo
[Scope]
RuntimeRandomizedExtraSec=24w

