# systemdlint/tests/auto/InvalidValue/2.50/good/Service_StartupAllowedCPUs_abc.target
[Unit]
Description=Foo
[Service]
StartupAllowedCPUs=abc

