# systemdlint/tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_9min.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=9min

