# systemdlint/tests/auto/InvalidValue/2.50/good/Service_AllowedMemoryNodes_abc.target
[Unit]
Description=Foo
[Service]
AllowedMemoryNodes=abc

