# systemdlint/tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_26year.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=26year

