# systemdlint/tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_5hr.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=5hr

