# systemdlint/tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_7minutes.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=7minutes

