# systemdlint/tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_3hours.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=3hours

