# systemdlint/tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_18sec.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=18sec

