# systemdlint/tests/auto/InvalidValue/2.50/good/Service_StartupAllowedMemoryNodes_abc.target
[Unit]
Description=Foo
[Service]
StartupAllowedMemoryNodes=abc

