# systemdlint/tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_15ms.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=15ms

