# systemdlint/tests/auto/InvalidValue/2.50/good/Service_RestrictNetworkInterfaces_abc.target
[Unit]
Description=Foo
[Service]
RestrictNetworkInterfaces=abc

