# systemdlint/tests/auto/InvalidValue/2.50/good/Service_TTYRows_9999999999999.target
[Unit]
Description=Foo
[Service]
TTYRows=9999999999999

