# systemdlint/tests/auto/InvalidValue/2.50/good/Service_ExitType_cgroup.target
[Unit]
Description=Foo
[Service]
ExitType=cgroup

