# systemdlint/tests/auto/InvalidValue/2.50/good/Service_RuntimeRandomizedExtraSec_1day.target
[Unit]
Description=Foo
[Service]
RuntimeRandomizedExtraSec=1day

