# systemdlint/tests/auto/InvalidValue/2.50/good/Slice_AllowedCPUs_abc.target
[Unit]
Description=Foo
[Slice]
AllowedCPUs=abc

