# systemdlint/tests/auto/InvalidValue/2.50/good/Socket_TTYColumns_0.target
[Unit]
Description=Foo
[Socket]
TTYColumns=0

