# systemdlint/tests/auto/InvalidValue/2.50/good/Socket_TTYRows_0.target
[Unit]
Description=Foo
[Socket]
TTYRows=0

