# systemdlint/tests/auto/InvalidValue/2.50/good/Socket_AllowedMemoryNodes_abc.target
[Unit]
Description=Foo
[Socket]
AllowedMemoryNodes=abc

