# systemdlint/tests/auto/InvalidValue/2.50/good/Socket_ExecSearchPath_abc.target
[Unit]
Description=Foo
[Socket]
ExecSearchPath=abc

