# systemdlint/tests/auto/InvalidValue/2.50/good/Socket_TTYColumns_9999999999999.target
[Unit]
Description=Foo
[Socket]
TTYColumns=9999999999999

