# systemdlint/tests/auto/InvalidValue/2.50/good/Swap_RestrictNetworkInterfaces_abc.target
[Unit]
Description=Foo
[Swap]
RestrictNetworkInterfaces=abc

