# systemdlint/tests/auto/InvalidValue/2.50/good/Swap_ExecSearchPath_abc.target
[Unit]
Description=Foo
[Swap]
ExecSearchPath=abc

