# systemdlint/tests/auto/InvalidValue/2.50/good/Swap_TTYColumns_9999999999999.target
[Unit]
Description=Foo
[Swap]
TTYColumns=9999999999999

