# systemdlint/tests/auto/InvalidValue/2.50/good/Swap_TTYRows_0.target
[Unit]
Description=Foo
[Swap]
TTYRows=0

