# systemdlint/tests/auto/InvalidValue/2.50/good/Swap_TTYRows_9999999999999.target
[Unit]
Description=Foo
[Swap]
TTYRows=9999999999999

