# systemdlint/tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_0days.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=0days

