# systemdlint/tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_22weeks.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=22weeks

