# systemdlint/tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_9min.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=9min

