# systemdlint/tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_19s.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=19s

