# systemdlint/tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_4hour.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=4hour

