# systemdlint/tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_13months.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=13months

