# systemdlint/tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_11m.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=11m

