# systemdlint/tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_8minute.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=8minute

