# systemdlint/tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_infinity.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=infinity

