# systemdlint/tests/auto/InvalidValue/2.50/good/Time_SaveIntervalSec_28%C2%B5s.target
[Unit]
Description=Foo
[Time]
SaveIntervalSec=28µs

