# systemdlint/tests/auto/InvalidValue/2.50/good/Unit_ConditionMemoryPressure_abc.target
[Unit]
Description=Foo
Description=abc
ConditionMemoryPressure=abc

