# tests/auto/OptionTooNew/2.46/good/BFIFO_Parent_abc.target
[Unit]
Description=Foo
[BFIFO]
Parent=123

