# tests/auto/OptionTooNew/2.46/good/DeficitRoundRobinScheduler_Parent_abc.target
[Unit]
Description=Foo
[DeficitRoundRobinScheduler]
Parent=123

