# tests/auto/OptionTooNew/2.46/good/DHCPServer_LRP_abc.target
[Unit]
Description=Foo
[DHCPServer]
LRP=123

