# tests/auto/OptionTooNew/2.46/good/DHCPv6_RequestOptions_0%209999999999999.target
[Unit]
Description=Foo
[DHCPv6]
RequestOptions=123

