# tests/auto/OptionTooNew/2.46/good/DHCPv6_WithoutRA_solicit.target
[Unit]
Description=Foo
[DHCPv6]
WithoutRA=123

