# tests/auto/OptionTooNew/2.46/good/EnhancedTransmissionSelection_Parent_abc.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
Parent=123

