# tests/auto/OptionTooNew/2.46/good/GenericRandomEarlyDetection_Parent_abc.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
Parent=123

