# tests/auto/OptionTooNew/2.46/good/GenericRandomEarlyDetection_VirtualQueues_1.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
VirtualQueues=123

