# tests/auto/OptionTooNew/2.46/good/HeavyHitterFilter_Parent_abc.target
[Unit]
Description=Foo
[HeavyHitterFilter]
Parent=123

