# tests/auto/OptionTooNew/2.46/good/HierarchyTokenBucketClass_CeilRate_0K.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
CeilRate=123

