# tests/auto/OptionTooNew/2.46/good/HierarchyTokenBucketClass_Rate_0K.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
Rate=123

