# tests/auto/OptionTooNew/2.46/good/HierarchyTokenBucketClass_OverheadBytes_0K.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
OverheadBytes=123

