# tests/auto/OptionTooNew/2.46/good/PFIFOHeadDrop_Parent_abc.target
[Unit]
Description=Foo
[PFIFOHeadDrop]
Parent=123

