# tests/auto/OptionTooNew/2.46/good/PIE_PacketLimit_0.target
[Unit]
Description=Foo
[PIE]
PacketLimit=123

