# tests/auto/OptionTooNew/2.46/good/PIE_Handle_0x1.target
[Unit]
Description=Foo
[PIE]
Handle=123

