# tests/auto/OptionTooNew/2.46/good/QuickFairQueueingClass_MaxPacketBytes_0K.target
[Unit]
Description=Foo
[QuickFairQueueingClass]
MaxPacketBytes=123

