# tests/auto/OptionTooNew/2.46/good/Service_RootHashSignature_abc.target
[Unit]
Description=Foo
[Service]
RootHashSignature=123

