# tests/auto/OptionTooNew/2.46/good/Socket_PassPacketInfo_yes.target
[Unit]
Description=Foo
[Socket]
PassPacketInfo=123

