# tests/auto/InvalidValue/2.46/bad/BFIFO_PacketLimit_-1.target
[Unit]
Description=Foo
[BFIFO]
PacketLimit=-1

