# tests/auto/InvalidValue/2.46/bad/BFIFO_PacketLimit_4294967295.target
[Unit]
Description=Foo
[BFIFO]
PacketLimit=4294967295

