# tests/auto/InvalidValue/2.46/bad/CAKE_Handle_0x0.target
[Unit]
Description=Foo
[CAKE]
Handle=0x0

