# tests/auto/InvalidValue/2.46/bad/CAKE_Bandwidth_10000000000000G.target
[Unit]
Description=Foo
[CAKE]
Bandwidth=10000000000000G

