# tests/auto/InvalidValue/2.46/bad/CAN_FDMode_baz.target
[Unit]
Description=Foo
[CAN]
FDMode=baz

