# tests/auto/InvalidValue/2.46/bad/CAN_ListenOnly_foo.target
[Unit]
Description=Foo
[CAN]
ListenOnly=foo

