# tests/auto/InvalidValue/2.46/bad/CAN_Termination_foo.target
[Unit]
Description=Foo
[CAN]
Termination=foo

