# tests/auto/InvalidValue/2.46/bad/CAN_BitRate_0M.target
[Unit]
Description=Foo
[CAN]
BitRate=0M

