# tests/auto/InvalidValue/2.46/bad/CAN_ListenOnly_baz.target
[Unit]
Description=Foo
[CAN]
ListenOnly=baz

