# tests/auto/InvalidValue/2.46/bad/CAN_DataBitRate_0K.target
[Unit]
Description=Foo
[CAN]
DataBitRate=0K

