# tests/auto/InvalidValue/2.46/bad/CAN_Termination_baz.target
[Unit]
Description=Foo
[CAN]
Termination=baz

