# tests/auto/InvalidValue/2.46/bad/CAN_FDMode_bar.target
[Unit]
Description=Foo
[CAN]
FDMode=bar

