# tests/auto/InvalidValue/2.46/bad/DHCPv6_RequestOptions_-1%2010000000000000.target
[Unit]
Description=Foo
[DHCPv6]
RequestOptions=-1 10000000000000

