# tests/auto/InvalidValue/2.46/bad/DHCPv6PrefixDelegation_Assign_baz.target
[Unit]
Description=Foo
[DHCPv6PrefixDelegation]
Assign=baz

