# tests/auto/InvalidValue/2.46/bad/DHCPServer_EmitLPR_bar.target
[Unit]
Description=Foo
[DHCPServer]
EmitLPR=bar

