# tests/auto/InvalidValue/2.46/bad/DHCPServer_EmitSMTP_foo.target
[Unit]
Description=Foo
[DHCPServer]
EmitSMTP=foo

