# tests/auto/InvalidValue/2.46/bad/EnhancedTransmissionSelection_QuantumBytes_10000000000000K.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
QuantumBytes=10000000000000K

