# tests/auto/InvalidValue/2.46/bad/EnhancedTransmissionSelection_PriorityMap_-1%2010000000000000.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
PriorityMap=-1 10000000000000

