# tests/auto/InvalidValue/2.46/bad/EnhancedTransmissionSelection_QuantumBytes_-1M.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
QuantumBytes=-1M

