# tests/auto/InvalidValue/2.46/bad/EnhancedTransmissionSelection_Handle_0x10000.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
Handle=0x10000

