# tests/auto/InvalidValue/2.46/bad/EnhancedTransmissionSelection_QuantumBytes_-1G.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
QuantumBytes=-1G

