# tests/auto/InvalidValue/2.46/bad/GenericRandomEarlyDetection_DefaultVirtualQueue_17.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
DefaultVirtualQueue=17

