# tests/auto/InvalidValue/2.46/bad/GenericRandomEarlyDetection_Handle_0x10000.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
Handle=0x10000

