# tests/auto/InvalidValue/2.46/bad/GenericRandomEarlyDetection_GenericRIO_foo.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
GenericRIO=foo

