# tests/auto/InvalidValue/2.46/bad/GenericRandomEarlyDetection_DefaultVirtualQueue_0.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
DefaultVirtualQueue=0

