# tests/auto/InvalidValue/2.46/bad/HeavyHitterFilter_Handle_0x10000.target
[Unit]
Description=Foo
[HeavyHitterFilter]
Handle=0x10000

