# tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_MTUBytes_-1M.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
MTUBytes=-1M

