# tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_Priority_-1.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
Priority=-1

