# tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_CeilRate_-1K.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
CeilRate=-1K

