# tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_QuantumBytes_-1K.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
QuantumBytes=-1K

