# tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_CeilRate_10000000000000K.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
CeilRate=10000000000000K

