# tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_Rate_-1K.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
Rate=-1K

