# tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_OverheadBytes_-1M.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
OverheadBytes=-1M

