# tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_CeilBufferBytes_-1G.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
CeilBufferBytes=-1G

