# tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_CeilBufferBytes_10000000000000K.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
CeilBufferBytes=10000000000000K

