# tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucket_DefaultClass_0x9184e72a000.target
[Unit]
Description=Foo
[HierarchyTokenBucket]
DefaultClass=0x9184e72a000

