# tests/auto/InvalidValue/2.46/bad/HierarchyTokenBucketClass_OverheadBytes_-1K.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
OverheadBytes=-1K

