# tests/auto/InvalidValue/2.46/bad/Link_AutoNegotiationFlowControl_baz.target
[Unit]
Description=Foo
[Link]
AutoNegotiationFlowControl=baz

