# tests/auto/InvalidValue/2.46/bad/Link_RxFlowControl_baz.target
[Unit]
Description=Foo
[Link]
RxFlowControl=baz

