# tests/auto/InvalidValue/2.46/bad/Link_TxFlowControl_bar.target
[Unit]
Description=Foo
[Link]
TxFlowControl=bar

