# tests/auto/InvalidValue/2.46/bad/Link_RxFlowControl_foo.target
[Unit]
Description=Foo
[Link]
RxFlowControl=foo

