# tests/auto/InvalidValue/2.46/bad/Link_AutoNegotiationFlowControl_bar.target
[Unit]
Description=Foo
[Link]
AutoNegotiationFlowControl=bar

