# tests/auto/InvalidValue/2.46/bad/Link_TxFlowControl_baz.target
[Unit]
Description=Foo
[Link]
TxFlowControl=baz

