# tests/auto/InvalidValue/2.46/bad/PFIFOFast_Handle_0x10000.target
[Unit]
Description=Foo
[PFIFOFast]
Handle=0x10000

