# tests/auto/InvalidValue/2.46/bad/PFIFOHeadDrop_Handle_0x0.target
[Unit]
Description=Foo
[PFIFOHeadDrop]
Handle=0x0

