# tests/auto/InvalidValue/2.46/bad/PFIFO_PacketLimit_-1.target
[Unit]
Description=Foo
[PFIFO]
PacketLimit=-1

