# tests/auto/InvalidValue/2.46/bad/PIE_Handle_0x10000.target
[Unit]
Description=Foo
[PIE]
Handle=0x10000

