# tests/auto/InvalidValue/2.46/bad/PIE_PacketLimit_-1.target
[Unit]
Description=Foo
[PIE]
PacketLimit=-1

