# tests/auto/InvalidValue/2.46/bad/QuickFairQueueing_Handle_0x10000.target
[Unit]
Description=Foo
[QuickFairQueueing]
Handle=0x10000

