# tests/auto/InvalidValue/2.46/bad/Socket_PassPacketInfo_baz.target
[Unit]
Description=Foo
[Socket]
PassPacketInfo=baz

