# tests/auto/InvalidValue/2.46/bad/Socket_PassPacketInfo_foo.target
[Unit]
Description=Foo
[Socket]
PassPacketInfo=foo

