# tests/auto/InvalidValue/2.46/bad/SR-IOV_Trust_bar.target
[Unit]
Description=Foo
[SR-IOV]
Trust=bar

