# tests/auto/InvalidValue/2.46/bad/StochasticFairBlue_PacketLimit_-1.target
[Unit]
Description=Foo
[StochasticFairBlue]
PacketLimit=-1

