# tests/auto/InvalidValue/2.46/bad/Unit_AssertPathIsEncrypted_1.target
[Unit]
Description=Foo
Description=abc
AssertPathIsEncrypted=1

