# tests/auto/InvalidValue/2.46/bad/Unit_AssertEnvironment_True.target
[Unit]
Description=Foo
Description=abc
AssertEnvironment=True

