# tests/auto/InvalidValue/2.46/bad/Unit_AssertPathIsEncrypted_True.target
[Unit]
Description=Foo
Description=abc
AssertPathIsEncrypted=True

