# tests/auto/InvalidValue/2.46/bad/Unit_ConditionEnvironment_True.target
[Unit]
Description=Foo
Description=abc
ConditionEnvironment=True

