# systemdlint/tests/auto/OptionTooNew/2.43/good/DHCPv6_UserClass_abc.target
[Unit]
Description=Foo
[DHCPv6]
UserClass=123

