# systemdlint/tests/auto/OptionTooNew/2.43/good/DHCPv6_IAID_0.target
[Unit]
Description=Foo
[DHCPv6]
IAID=123

