# systemdlint/tests/auto/OptionTooNew/2.43/good/DHCPv4_RouteTable_1.target
[Unit]
Description=Foo
[DHCPv4]
RouteTable=123

