# systemdlint/tests/auto/OptionTooNew/2.43/good/FooOverUDP_Peer_1.1.1.1--0.target
[Unit]
Description=Foo
[FooOverUDP]
Peer=123

