# systemdlint/tests/auto/OptionTooNew/2.43/good/FooOverUDP_PeerPort_1.target
[Unit]
Description=Foo
[FooOverUDP]
PeerPort=123

