# systemdlint/tests/auto/OptionTooNew/2.43/good/Service_NUMAPolicy_default.target
[Unit]
Description=Foo
[Service]
NUMAPolicy=123

