# systemdlint/tests/auto/OptionTooNew/2.43/good/Swap_DefaultMemoryLow_0K.target
[Unit]
Description=Foo
[Swap]
DefaultMemoryLow=123

