# systemdlint/tests/auto/OptionTooNew/2.43/good/VXLAN_GenericProtocolExtension_yes.target
[Unit]
Description=Foo
[VXLAN]
GenericProtocolExtension=123

