# systemdlint/tests/auto/OptionTooNew/2.48/bad/Link_WakeOnLan_off%20phy%20unicast%20multicast%20broadcast%20arp%20magic%20secureon.target
[Unit]
Description=Foo
[Link]
WakeOnLan=123

