# systemdlint/tests/auto/OptionTooNew/2.48/bad/Manager_StatusUnitFormat_name.target
[Unit]
Description=Foo
[Manager]
StatusUnitFormat=123

