# systemdlint/tests/auto/OptionTooNew/2.48/bad/Swap_SocketBindAllow_abc.target
[Unit]
Description=Foo
[Swap]
SocketBindAllow=123

