# systemdlint/tests/auto/OptionTooNew/2.48/bad/Upload_NetworkTimeoutSec_0days.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=123

