# systemdlint/tests/auto/OptionTooNew/2.48/bad/Upload_ServerCertificateFile_abc.target
[Unit]
Description=Foo
[Upload]
ServerCertificateFile=123

