# systemdlint/tests/auto/OptionTooNew/2.48/bad/Upload_URL_1.1.1.1--0.target
[Unit]
Description=Foo
[Upload]
URL=123

