# tests/auto/OptionTooNew/2.31/good/Service_InaccessableDirectories_abc.target
[Unit]
Description=Foo
[Service]
InaccessableDirectories=123

