# tests/auto/OptionTooNew/2.31/good/Service_RestrictRealtime_yes.target
[Unit]
Description=Foo
[Service]
RestrictRealtime=123

