# tests/auto/OptionTooNew/2.31/good/Socket_ReadOnlyDirectories_abc.target
[Unit]
Description=Foo
[Socket]
ReadOnlyDirectories=123

