# tests/auto/InvalidValue/2.30/good/Mount_Capabilities_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
Capabilities=abc

