# tests/auto/InvalidValue/2.30/good/Service_Capabilities_abc.target
[Unit]
Description=Foo
[Service]
Capabilities=abc

