# tests/auto/InvalidValue/2.30/good/Socket_Capabilities_abc.target
[Unit]
Description=Foo
[Socket]
Capabilities=abc

