# tests/auto/InvalidValue/2.30/good/Swap_Capabilities_abc.target
[Unit]
Description=Foo
[Swap]
Capabilities=abc

