# systemdlint/tests/auto/OptionTooNew/2.11/bad/Unit_ConditionSecurity_abc.target
[Unit]
Description=Foo
Description=abc
ConditionSecurity=123

