# systemdlint/tests/auto/OptionTooNew/2.45/bad/BFIFO_PacketLimit_0.target
[Unit]
Description=Foo
[BFIFO]
PacketLimit=123

