# systemdlint/tests/auto/OptionTooNew/2.45/bad/CAKE_Parent_abc.target
[Unit]
Description=Foo
[CAKE]
Parent=123

