# systemdlint/tests/auto/OptionTooNew/2.45/bad/CAN_ListenOnly_yes.target
[Unit]
Description=Foo
[CAN]
ListenOnly=123

