# systemdlint/tests/auto/OptionTooNew/2.45/bad/CAN_Termination_yes.target
[Unit]
Description=Foo
[CAN]
Termination=123

