# systemdlint/tests/auto/OptionTooNew/2.45/bad/DHCPServer_SMTP_abc.target
[Unit]
Description=Foo
[DHCPServer]
SMTP=123

