# systemdlint/tests/auto/OptionTooNew/2.45/bad/DHCPv6_VendorClass_abc.target
[Unit]
Description=Foo
[DHCPv6]
VendorClass=123

