# systemdlint/tests/auto/OptionTooNew/2.45/bad/EnhancedTransmissionSelection_PriorityMap_abc.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
PriorityMap=123

