# systemdlint/tests/auto/OptionTooNew/2.45/bad/HierarchyTokenBucket_RateToQuantum_0.target
[Unit]
Description=Foo
[HierarchyTokenBucket]
RateToQuantum=123

