# systemdlint/tests/auto/OptionTooNew/2.45/bad/HierarchyTokenBucketClass_ClassId_abc.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
ClassId=123

