# systemdlint/tests/auto/OptionTooNew/2.45/bad/Link_Group_0.target
[Unit]
Description=Foo
[Link]
Group=123

