# systemdlint/tests/auto/OptionTooNew/2.45/bad/PIE_Parent_abc.target
[Unit]
Description=Foo
[PIE]
Parent=123

