# systemdlint/tests/auto/OptionTooNew/2.45/bad/QuickFairQueueingClass_Weight_1.target
[Unit]
Description=Foo
[QuickFairQueueingClass]
Weight=123

