# systemdlint/tests/auto/OptionTooNew/2.45/bad/QuickFairQueueingClass_Parent_abc.target
[Unit]
Description=Foo
[QuickFairQueueingClass]
Parent=123

