# systemdlint/tests/auto/OptionTooNew/2.45/bad/Swap_RootHashSignature_abc.target
[Unit]
Description=Foo
[Swap]
RootHashSignature=123

