# systemdlint/tests/auto/OptionTooNew/2.45/bad/Unit_AssertEnvironment_a%3Db.target
[Unit]
Description=Foo
Description=abc
AssertEnvironment=123

