# systemdlint/tests/auto/OptionTooNew/2.45/bad/Unit_ConditionEnvironment_a%3Db.target
[Unit]
Description=Foo
Description=abc
ConditionEnvironment=123

