# tests/auto/OptionTooNew/2.35/good/Service_RuntimeDirectory_abc.target
[Unit]
Description=Foo
[Service]
RuntimeDirectory=123

