# tests/auto/OptionTooNew/2.35/good/Service_CacheDirectory_abc.target
[Unit]
Description=Foo
[Service]
CacheDirectory=123

