# tests/auto/OptionTooNew/2.35/good/Service_StateDirectory_abc.target
[Unit]
Description=Foo
[Service]
StateDirectory=123

