# tests/auto/OptionTooNew/2.35/good/Service_UnsetEnvironment_abc.target
[Unit]
Description=Foo
[Service]
UnsetEnvironment=123

