# tests/auto/OptionTooNew/2.35/good/Socket_CacheDirectory_abc.target
[Unit]
Description=Foo
[Socket]
CacheDirectory=123

