# tests/auto/OptionTooNew/2.35/good/Socket_IPAddressDeny_1.1.1.1--0.target
[Unit]
Description=Foo
[Socket]
IPAddressDeny=123

