# tests/auto/OptionTooNew/2.35/good/Socket_LogsDirectory_abc.target
[Unit]
Description=Foo
[Socket]
LogsDirectory=123

