# tests/auto/OptionTooNew/2.35/good/Swap_RuntimeDirectory_abc.target
[Unit]
Description=Foo
[Swap]
RuntimeDirectory=123

