# tests/auto/OptionTooNew/2.35/good/Swap_StateDirectory_abc.target
[Unit]
Description=Foo
[Swap]
StateDirectory=123

