# tests/auto/OptionTooNew/2.35/good/Swap_LockPersonality_yes.target
[Unit]
Description=Foo
[Swap]
LockPersonality=123

