# tests/auto/OptionTooNew/2.35/good/Swap_UnsetEnvironment_abc.target
[Unit]
Description=Foo
[Swap]
UnsetEnvironment=123

