# tests/auto/OptionTooNew/2.35/good/Swap_IPAccounting_yes.target
[Unit]
Description=Foo
[Swap]
IPAccounting=123

