# systemdlint/tests/manual/NoFailureCheck/2.31/good/Service_ExecStart_bin-busybox.service
[Unit]
Description=Foo

[Service]
Type=oneshot
ExecStart=/bin/busybox