# systemdlint/tests/manual/NoFailureCheck/2.31/good/Service_ExecStart_ls@.service
[Unit]
Description=Foo

[Service]
Type=oneshot
ExecStartPre=/bin/ls /foo/bar/%i
ExecStart=/bin/ls -l /foo/bar/%i
ExecStartPost=/bin/ls /foo/bar-zap/%i
