# tests/auto/OptionTooNew/2.18/good/Socket_SmackProcessLabel_abc.target
[Unit]
Description=Foo
[Socket]
SmackProcessLabel=123

