# tests/auto/OptionTooNew/2.18/good/Swap_SmackProcessLabel_abc.target
[Unit]
Description=Foo
[Swap]
SmackProcessLabel=123

