# systemdlint/tests/auto/OptionTooNew/2.41/bad/Socket_ProtectHostname_yes.target
[Unit]
Description=Foo
[Socket]
ProtectHostname=123

