# systemdlint/tests/auto/OptionTooNew/2.43/bad/DHCPServer_SIP_abc.target
[Unit]
Description=Foo
[DHCPServer]
SIP=123

