# systemdlint/tests/auto/OptionTooNew/2.43/bad/DHCPServer_SendOption_abc.target
[Unit]
Description=Foo
[DHCPServer]
SendOption=123

