# systemdlint/tests/auto/OptionTooNew/2.43/bad/Swap_ProtectKernelLogs_yes.target
[Unit]
Description=Foo
[Swap]
ProtectKernelLogs=123

