# systemdlint/tests/auto/OptionTooNew/2.43/bad/TrafficControlQueueingDiscipline_NetworkEmulatorDelayJitterSec_0.target
[Unit]
Description=Foo
[TrafficControlQueueingDiscipline]
NetworkEmulatorDelayJitterSec=123

