# systemdlint/tests/auto/OptionTooNew/2.43/bad/TrafficControlQueueingDiscipline_NetworkEmulatorDuplicateRate_0%25.target
[Unit]
Description=Foo
[TrafficControlQueueingDiscipline]
NetworkEmulatorDuplicateRate=123

