# tests/auto/OptionTooNew/2.48/good/Link_ActivationPolicy_up.target
[Unit]
Description=Foo
[Link]
ActivationPolicy=123

