# tests/auto/OptionTooNew/2.48/good/RoutingPolicyRule_Family_ipv4.target
[Unit]
Description=Foo
[RoutingPolicyRule]
Family=123

