# tests/auto/OptionTooNew/2.48/good/Socket_ManagedOOMPreference_none.target
[Unit]
Description=Foo
[Socket]
ManagedOOMPreference=123

