# tests/auto/OptionTooNew/2.48/good/Socket_NoExecPaths_abc.target
[Unit]
Description=Foo
[Socket]
NoExecPaths=123

