# tests/auto/OptionTooNew/2.48/good/VLAN_IngressQOSMaps_abc.target
[Unit]
Description=Foo
[VLAN]
IngressQOSMaps=123

