# tests/auto/OptionTooNew/2.46/bad/IPv6SendRA_DNS_abc.target
[Unit]
Description=Foo
[IPv6SendRA]
DNS=123

