# tests/auto/OptionTooNew/2.46/bad/IPv6SendRA_EmitDNS_yes.target
[Unit]
Description=Foo
[IPv6SendRA]
EmitDNS=123

