# tests/auto/OptionTooNew/2.46/bad/IPv6SendRA_Domains_abc.target
[Unit]
Description=Foo
[IPv6SendRA]
Domains=123

