# tests/auto/OptionTooNew/2.46/bad/Mount_RootImageOptions_%20root%20root-secondary%20home%20srv%20esp%20xbootldr%20tmp%20var%20usr.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
RootImageOptions=123

