# tests/auto/OptionTooNew/2.46/bad/Service_RootImageOptions_%20root%20root-secondary%20home%20srv%20esp%20xbootldr%20tmp%20var%20usr.target
[Unit]
Description=Foo
[Service]
RootImageOptions=123

