# tests/auto/OptionTooNew/2.46/bad/Service_MountImages_abc.target
[Unit]
Description=Foo
[Service]
MountImages=123

