# tests/auto/OptionTooNew/2.46/bad/Socket_Timestamping_off.target
[Unit]
Description=Foo
[Socket]
Timestamping=123

