# tests/auto/OptionTooNew/2.46/bad/Socket_ManagedOOMSwap_abc.target
[Unit]
Description=Foo
[Socket]
ManagedOOMSwap=123

