# tests/auto/OptionTooNew/2.46/bad/Swap_ProtectProc_noaccess.target
[Unit]
Description=Foo
[Swap]
ProtectProc=123

