# tests/auto/OptionTooNew/2.46/bad/Swap_NUMAMask_all.target
[Unit]
Description=Foo
[Swap]
NUMAMask=123

