# tests/auto/OptionTooNew/2.46/bad/Swap_MountImages_abc.target
[Unit]
Description=Foo
[Swap]
MountImages=123

