# systemdlint/tests/auto/OptionTooNew/2.11/good/Socket_RestrictAddressFamilies_abc.target
[Unit]
Description=Foo
[Socket]
RestrictAddressFamilies=123

