# systemdlint/tests/auto/OptionTooNew/2.11/good/Swap_RestrictAddressFamilies_abc.target
[Unit]
Description=Foo
[Swap]
RestrictAddressFamilies=123

