# tests/auto/InvalidValue/2.32/good/Mount_ProtectKernelModules_1.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
ProtectKernelModules=1

