# tests/auto/InvalidValue/2.32/good/Service_ProtectKernelTunables_on.target
[Unit]
Description=Foo
[Service]
ProtectKernelTunables=on

