# tests/auto/InvalidValue/2.32/good/Service_DynamicUser_false.target
[Unit]
Description=Foo
[Service]
DynamicUser=false

