# tests/auto/InvalidValue/2.32/good/Service_ProtectKernelTunables_0.target
[Unit]
Description=Foo
[Service]
ProtectKernelTunables=0

