# tests/auto/InvalidValue/2.32/good/Service_DynamicUser_true.target
[Unit]
Description=Foo
[Service]
DynamicUser=true

