# tests/auto/InvalidValue/2.32/good/Service_CPUAccounting_false.target
[Unit]
Description=Foo
[Service]
CPUAccounting=false

