# tests/auto/InvalidValue/2.32/good/Service_ProtectKernelTunables_1.target
[Unit]
Description=Foo
[Service]
ProtectKernelTunables=1

