# tests/auto/InvalidValue/2.32/good/Service_DynamicUser_yes.target
[Unit]
Description=Foo
[Service]
DynamicUser=yes

