# tests/auto/InvalidValue/2.32/good/Service_ProtectControlGroups_1.target
[Unit]
Description=Foo
[Service]
ProtectControlGroups=1

