# tests/auto/InvalidValue/2.32/good/Service_CPUAccounting_1.target
[Unit]
Description=Foo
[Service]
CPUAccounting=1

