# tests/auto/InvalidValue/2.32/good/Service_ProtectKernelTunables_yes.target
[Unit]
Description=Foo
[Service]
ProtectKernelTunables=yes

