# tests/auto/InvalidValue/2.32/good/Service_CPUWeight_1.target
[Unit]
Description=Foo
[Service]
CPUWeight=1

