# tests/auto/InvalidValue/2.32/good/Service_PrivateUsers_yes.target
[Unit]
Description=Foo
[Service]
PrivateUsers=yes

