# tests/auto/InvalidValue/2.32/good/Service_PrivateUsers_false.target
[Unit]
Description=Foo
[Service]
PrivateUsers=false

