# tests/auto/InvalidValue/2.32/good/Service_PrivateUsers_no.target
[Unit]
Description=Foo
[Service]
PrivateUsers=no

