# tests/auto/InvalidValue/2.32/good/Service_DynamicUser_no.target
[Unit]
Description=Foo
[Service]
DynamicUser=no

