# tests/auto/InvalidValue/2.32/good/Service_ProtectKernelModules_no.target
[Unit]
Description=Foo
[Service]
ProtectKernelModules=no

