# tests/auto/InvalidValue/2.32/good/Service_ProtectControlGroups_true.target
[Unit]
Description=Foo
[Service]
ProtectControlGroups=true

