# tests/auto/InvalidValue/2.32/good/Service_ProtectKernelModules_true.target
[Unit]
Description=Foo
[Service]
ProtectKernelModules=true

