# tests/auto/InvalidValue/2.32/good/Service_DynamicUser_off.target
[Unit]
Description=Foo
[Service]
DynamicUser=off

