# tests/auto/InvalidValue/2.32/good/Service_CPUAccounting_on.target
[Unit]
Description=Foo
[Service]
CPUAccounting=on

