# tests/auto/InvalidValue/2.32/good/Service_ProtectKernelTunables_no.target
[Unit]
Description=Foo
[Service]
ProtectKernelTunables=no

