# tests/auto/InvalidValue/2.32/good/Service_ProtectKernelModules_yes.target
[Unit]
Description=Foo
[Service]
ProtectKernelModules=yes

