# tests/auto/InvalidValue/2.32/good/Service_PrivateUsers_0.target
[Unit]
Description=Foo
[Service]
PrivateUsers=0

