# tests/auto/InvalidValue/2.32/good/Socket_CPUAccounting_off.target
[Unit]
Description=Foo
[Socket]
CPUAccounting=off

