# tests/auto/InvalidValue/2.32/good/Socket_ProtectKernelTunables_true.target
[Unit]
Description=Foo
[Socket]
ProtectKernelTunables=true

