# tests/auto/InvalidValue/2.32/good/Socket_ProtectControlGroups_off.target
[Unit]
Description=Foo
[Socket]
ProtectControlGroups=off

