# tests/auto/InvalidValue/2.32/good/Socket_DynamicUser_false.target
[Unit]
Description=Foo
[Socket]
DynamicUser=false

