# tests/auto/InvalidValue/2.32/good/Socket_PrivateUsers_on.target
[Unit]
Description=Foo
[Socket]
PrivateUsers=on

