# tests/auto/InvalidValue/2.32/good/Socket_ProtectControlGroups_true.target
[Unit]
Description=Foo
[Socket]
ProtectControlGroups=true

