# tests/auto/InvalidValue/2.32/good/Socket_MaxConnections_0.target
[Unit]
Description=Foo
[Socket]
MaxConnections=0

