# tests/auto/InvalidValue/2.32/good/Socket_CPUAccounting_no.target
[Unit]
Description=Foo
[Socket]
CPUAccounting=no

