# tests/auto/InvalidValue/2.32/good/Socket_CPUAccounting_1.target
[Unit]
Description=Foo
[Socket]
CPUAccounting=1

