# tests/auto/InvalidValue/2.32/good/Socket_DynamicUser_yes.target
[Unit]
Description=Foo
[Socket]
DynamicUser=yes

