# tests/auto/InvalidValue/2.32/good/Socket_PrivateUsers_1.target
[Unit]
Description=Foo
[Socket]
PrivateUsers=1

