# tests/auto/InvalidValue/2.32/good/Socket_ProtectControlGroups_no.target
[Unit]
Description=Foo
[Socket]
ProtectControlGroups=no

