# tests/auto/InvalidValue/2.32/good/Socket_ProtectKernelTunables_on.target
[Unit]
Description=Foo
[Socket]
ProtectKernelTunables=on

