# tests/auto/InvalidValue/2.32/good/Socket_DynamicUser_no.target
[Unit]
Description=Foo
[Socket]
DynamicUser=no

