# tests/auto/InvalidValue/2.32/good/Socket_MemorySwapMax_9999999999999B.target
[Unit]
Description=Foo
[Socket]
MemorySwapMax=9999999999999B

