# tests/auto/InvalidValue/2.32/good/Socket_DynamicUser_on.target
[Unit]
Description=Foo
[Socket]
DynamicUser=on

