# tests/auto/InvalidValue/2.32/good/Socket_ProtectControlGroups_0.target
[Unit]
Description=Foo
[Socket]
ProtectControlGroups=0

