# tests/auto/InvalidValue/2.32/good/Socket_ProtectControlGroups_on.target
[Unit]
Description=Foo
[Socket]
ProtectControlGroups=on

