# tests/auto/InvalidValue/2.32/good/Socket_ProtectControlGroups_1.target
[Unit]
Description=Foo
[Socket]
ProtectControlGroups=1

