# tests/auto/InvalidValue/2.32/good/Socket_ProtectKernelTunables_yes.target
[Unit]
Description=Foo
[Socket]
ProtectKernelTunables=yes

