# tests/auto/InvalidValue/2.32/good/Socket_DynamicUser_off.target
[Unit]
Description=Foo
[Socket]
DynamicUser=off

