# tests/auto/InvalidValue/2.32/good/Socket_ProtectKernelModules_on.target
[Unit]
Description=Foo
[Socket]
ProtectKernelModules=on

