# tests/auto/InvalidValue/2.32/good/Socket_ProtectKernelModules_off.target
[Unit]
Description=Foo
[Socket]
ProtectKernelModules=off

