# tests/auto/InvalidValue/2.32/good/Socket_MaxConnections_9999999999999.target
[Unit]
Description=Foo
[Socket]
MaxConnections=9999999999999

