# tests/auto/InvalidValue/2.32/good/Socket_PrivateUsers_true.target
[Unit]
Description=Foo
[Socket]
PrivateUsers=true

