# tests/auto/InvalidValue/2.32/good/Socket_CPUAccounting_on.target
[Unit]
Description=Foo
[Socket]
CPUAccounting=on

