# tests/auto/InvalidValue/2.32/good/Socket_ProtectControlGroups_false.target
[Unit]
Description=Foo
[Socket]
ProtectControlGroups=false

