# tests/auto/InvalidValue/2.32/good/Socket_PrivateUsers_no.target
[Unit]
Description=Foo
[Socket]
PrivateUsers=no

