# tests/auto/InvalidValue/2.32/good/Socket_ProtectKernelTunables_no.target
[Unit]
Description=Foo
[Socket]
ProtectKernelTunables=no

