# tests/auto/InvalidValue/2.32/good/Swap_PrivateUsers_0.target
[Unit]
Description=Foo
[Swap]
PrivateUsers=0

