# tests/auto/InvalidValue/2.32/good/Swap_CPUAccounting_on.target
[Unit]
Description=Foo
[Swap]
CPUAccounting=on

