# tests/auto/InvalidValue/2.32/good/Swap_ProtectKernelModules_false.target
[Unit]
Description=Foo
[Swap]
ProtectKernelModules=false

