# tests/auto/InvalidValue/2.32/good/Swap_CPUAccounting_off.target
[Unit]
Description=Foo
[Swap]
CPUAccounting=off

