# tests/auto/InvalidValue/2.32/good/Swap_ProtectKernelTunables_yes.target
[Unit]
Description=Foo
[Swap]
ProtectKernelTunables=yes

