# tests/auto/InvalidValue/2.32/good/Swap_CPUAccounting_yes.target
[Unit]
Description=Foo
[Swap]
CPUAccounting=yes

