# tests/auto/InvalidValue/2.32/good/Swap_DynamicUser_off.target
[Unit]
Description=Foo
[Swap]
DynamicUser=off

