# tests/auto/InvalidValue/2.32/good/Swap_CPUWeight_10000.target
[Unit]
Description=Foo
[Swap]
CPUWeight=10000

