# tests/auto/InvalidValue/2.32/good/Swap_ProtectKernelTunables_no.target
[Unit]
Description=Foo
[Swap]
ProtectKernelTunables=no

