# tests/auto/InvalidValue/2.32/good/Swap_DynamicUser_yes.target
[Unit]
Description=Foo
[Swap]
DynamicUser=yes

