# tests/auto/InvalidValue/2.32/good/Swap_DynamicUser_false.target
[Unit]
Description=Foo
[Swap]
DynamicUser=false

