# tests/auto/InvalidValue/2.32/good/Swap_DynamicUser_0.target
[Unit]
Description=Foo
[Swap]
DynamicUser=0

