# tests/auto/InvalidValue/2.32/good/Swap_ProtectKernelModules_yes.target
[Unit]
Description=Foo
[Swap]
ProtectKernelModules=yes

