# tests/auto/InvalidValue/2.32/good/Swap_CPUWeight_1.target
[Unit]
Description=Foo
[Swap]
CPUWeight=1

