# tests/auto/InvalidValue/2.32/good/Swap_PrivateUsers_1.target
[Unit]
Description=Foo
[Swap]
PrivateUsers=1

