# tests/auto/InvalidValue/2.32/good/Swap_CPUAccounting_no.target
[Unit]
Description=Foo
[Swap]
CPUAccounting=no

