# tests/auto/InvalidValue/2.32/good/Swap_ProtectControlGroups_true.target
[Unit]
Description=Foo
[Swap]
ProtectControlGroups=true

