# tests/auto/InvalidValue/2.32/good/Swap_ProtectKernelModules_off.target
[Unit]
Description=Foo
[Swap]
ProtectKernelModules=off

