# systemdlint/tests/auto/InvalidValue/2.34/good/Mount_StandardInputData_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
StandardInputData=abc

