# systemdlint/tests/auto/InvalidValue/2.34/good/Service_LogLevelMax_err.target
[Unit]
Description=Foo
[Service]
LogLevelMax=err

