# systemdlint/tests/auto/InvalidValue/2.34/good/Service_LogLevelMax_notice.target
[Unit]
Description=Foo
[Service]
LogLevelMax=notice

