# systemdlint/tests/auto/InvalidValue/2.34/good/Service_StandardInput_tty-fail.target
[Unit]
Description=Foo
[Service]
StandardInput=tty-fail

