# systemdlint/tests/auto/InvalidValue/2.34/good/Service_StandardInput_abc.target
[Unit]
Description=Foo
[Service]
StandardInput=abc

