# systemdlint/tests/auto/InvalidValue/2.34/good/Service_LogLevelMax_warning.target
[Unit]
Description=Foo
[Service]
LogLevelMax=warning

