# systemdlint/tests/auto/InvalidValue/2.34/good/Service_StandardInput_socket.target
[Unit]
Description=Foo
[Service]
StandardInput=socket

