# systemdlint/tests/auto/InvalidValue/2.34/good/Service_LogLevelMax_crit.target
[Unit]
Description=Foo
[Service]
LogLevelMax=crit

