# systemdlint/tests/auto/InvalidValue/2.34/good/Socket_StandardInput_tty.target
[Unit]
Description=Foo
[Socket]
StandardInput=tty

