# systemdlint/tests/auto/InvalidValue/2.34/good/Socket_LogLevelMax_warning.target
[Unit]
Description=Foo
[Socket]
LogLevelMax=warning

