# systemdlint/tests/auto/InvalidValue/2.34/good/Socket_StandardInput_abc.target
[Unit]
Description=Foo
[Socket]
StandardInput=abc

