# systemdlint/tests/auto/InvalidValue/2.34/good/Socket_LogLevelMax_err.target
[Unit]
Description=Foo
[Socket]
LogLevelMax=err

