# systemdlint/tests/auto/InvalidValue/2.34/good/Socket_LogLevelMax_crit.target
[Unit]
Description=Foo
[Socket]
LogLevelMax=crit

