# systemdlint/tests/auto/InvalidValue/2.34/good/Socket_LogLevelMax_debug.target
[Unit]
Description=Foo
[Socket]
LogLevelMax=debug

