# systemdlint/tests/auto/InvalidValue/2.34/good/Socket_StandardInputData_abc.target
[Unit]
Description=Foo
[Socket]
StandardInputData=abc

