# systemdlint/tests/auto/InvalidValue/2.34/good/Swap_StandardInput_null.target
[Unit]
Description=Foo
[Swap]
StandardInput=null

