# systemdlint/tests/auto/InvalidValue/2.34/good/Swap_LogLevelMax_debug.target
[Unit]
Description=Foo
[Swap]
LogLevelMax=debug

