# systemdlint/tests/auto/InvalidValue/2.34/good/Swap_StandardInput_tty.target
[Unit]
Description=Foo
[Swap]
StandardInput=tty

