# systemdlint/tests/auto/InvalidValue/2.34/good/Swap_LogLevelMax_crit.target
[Unit]
Description=Foo
[Swap]
LogLevelMax=crit

