# systemdlint/tests/auto/InvalidValue/2.34/good/Swap_LogLevelMax_emerg.target
[Unit]
Description=Foo
[Swap]
LogLevelMax=emerg

