# systemdlint/tests/auto/InvalidValue/2.34/good/Unit_CollectMode_inactive.target
[Unit]
Description=Foo
Description=abc
CollectMode=inactive

