# systemdlint/tests/auto/InvalidValue/2.34/good/Unit_SuccessAction_reboot-force.target
[Unit]
Description=Foo
Description=abc
SuccessAction=reboot-force

