# systemdlint/tests/auto/InvalidValue/2.34/good/Unit_ConditionGroup_abc.target
[Unit]
Description=Foo
Description=abc
ConditionGroup=abc

