# systemdlint/tests/auto/OptionTooNew/2.33/good/Mount_BindPaths_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
BindPaths=123

