# tests/auto/InvalidValue/2.49/good/Address_RouteMetric_0.target
[Unit]
Description=Foo
[Address]
Address=abc
RouteMetric=0

