# tests/auto/InvalidValue/2.49/good/DHCPServer_ServerAddress_255.255.255.255--32.target
[Unit]
Description=Foo
[DHCPServer]
ServerAddress=255.255.255.255/32

