# tests/auto/InvalidValue/2.49/good/DHCPv4_RoutesToNTP_off.target
[Unit]
Description=Foo
[DHCPv4]
RoutesToNTP=off

