# tests/auto/InvalidValue/2.49/good/DHCPv4_RoutesToNTP_true.target
[Unit]
Description=Foo
[DHCPv4]
RoutesToNTP=true

