# tests/auto/InvalidValue/2.49/good/DHCPServer_UplinkInterface_abc.target
[Unit]
Description=Foo
[DHCPServer]
UplinkInterface=abc

