# tests/auto/InvalidValue/2.49/good/Link_RequiredFamilyForOnline_any.target
[Unit]
Description=Foo
[Link]
RequiredFamilyForOnline=any

