# tests/auto/InvalidValue/2.49/good/Link_WakeOnLan_off%20phy%20unicast%20multicast%20broadcast%20arp%20magic%20secureon.target
[Unit]
Description=Foo
[Link]
WakeOnLan=off phy unicast multicast broadcast arp magic secureon

