# tests/auto/InvalidValue/2.49/good/Manager_StatusUnitFormat_combined.target
[Unit]
Description=Foo
[Manager]
StatusUnitFormat=combined

