# tests/auto/InvalidValue/2.49/good/Manager_StatusUnitFormat_description.target
[Unit]
Description=Foo
[Manager]
StatusUnitFormat=description

