# tests/auto/InvalidValue/2.49/good/Mount_SocketBindAllow_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
SocketBindAllow=abc

