# tests/auto/InvalidValue/2.49/good/Network_ManageForeignRoutingPolicyRules_off.target
[Unit]
Description=Foo
[Network]
ManageForeignRoutingPolicyRules=off

