# tests/auto/InvalidValue/2.49/good/Network_ManageForeignRoutingPolicyRules_1.target
[Unit]
Description=Foo
[Network]
ManageForeignRoutingPolicyRules=1

