# tests/auto/InvalidValue/2.49/good/Scope_SocketBindAllow_abc.target
[Unit]
Description=Foo
[Scope]
SocketBindAllow=abc

