# tests/auto/InvalidValue/2.49/good/Slice_SocketBindAllow_abc.target
[Unit]
Description=Foo
[Slice]
SocketBindAllow=abc

