# tests/auto/InvalidValue/2.49/good/Socket_SocketBindAllow_abc.target
[Unit]
Description=Foo
[Socket]
SocketBindAllow=abc

