# tests/auto/InvalidValue/2.49/good/Socket_SocketBindDeny_abc.target
[Unit]
Description=Foo
[Socket]
SocketBindDeny=abc

