# tests/auto/InvalidValue/2.49/good/Swap_SocketBindDeny_abc.target
[Unit]
Description=Foo
[Swap]
SocketBindDeny=abc

