# tests/auto/InvalidValue/2.49/good/Swap_SocketBindAllow_abc.target
[Unit]
Description=Foo
[Swap]
SocketBindAllow=abc

