# tests/auto/InvalidValue/2.49/good/Upload_NetworkTimeoutSec_11m.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=11m

