# tests/auto/InvalidValue/2.49/good/Upload_URL_http%3A----user%3Apwd%40test.com%3A1234.target
[Unit]
Description=Foo
[Upload]
URL=http://user:pwd@test.com:1234

