# tests/auto/InvalidValue/2.49/good/Upload_NetworkTimeoutSec_14msec.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=14msec

