# tests/auto/InvalidValue/2.49/good/Upload_NetworkTimeoutSec_2d.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=2d

