# tests/auto/InvalidValue/2.49/good/Upload_NetworkTimeoutSec_0days.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=0days

