# tests/auto/InvalidValue/2.49/good/Upload_NetworkTimeoutSec_5hr.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=5hr

