# tests/auto/InvalidValue/2.49/good/Upload_NetworkTimeoutSec_16seconds.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=16seconds

