# tests/auto/InvalidValue/2.49/good/Upload_NetworkTimeoutSec_25years.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=25years

