# tests/auto/InvalidValue/2.49/good/Upload_TrustedCertificateFile_abc.target
[Unit]
Description=Foo
[Upload]
TrustedCertificateFile=abc

