# tests/auto/InvalidValue/2.49/good/Upload_NetworkTimeoutSec_7minutes.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=7minutes

