# tests/auto/InvalidValue/2.49/good/Upload_URL_http%3A----a.b.com.target
[Unit]
Description=Foo
[Upload]
URL=http://a.b.com

