# tests/auto/InvalidValue/2.49/good/Upload_NetworkTimeoutSec_infinity.target
[Unit]
Description=Foo
[Upload]
NetworkTimeoutSec=infinity

