# tests/auto/InvalidValue/2.49/good/Upload_ServerCertificateFile_abc.target
[Unit]
Description=Foo
[Upload]
ServerCertificateFile=abc

