# tests/auto/InvalidValue/2.49/bad/Address_RouteMetric_4294967296.target
[Unit]
Description=Foo
[Address]
Address=abc
RouteMetric=4294967296

