# tests/auto/InvalidValue/2.49/bad/DHCPServerStaticLease_MACAddress_01%3A23%3A45%3A67%3A89.target
[Unit]
Description=Foo
[DHCPServerStaticLease]
MACAddress=01:23:45:67:89

