# tests/auto/InvalidValue/2.49/bad/DHCPServer_ServerAddress_0--abc.target
[Unit]
Description=Foo
[DHCPServer]
ServerAddress=0/abc

