# tests/auto/InvalidValue/2.49/bad/DHCPServerStaticLease_MACAddress_True.target
[Unit]
Description=Foo
[DHCPServerStaticLease]
MACAddress=True

