# tests/auto/InvalidValue/2.49/bad/DHCPServer_ServerAddress_abc.target
[Unit]
Description=Foo
[DHCPServer]
ServerAddress=abc

