# tests/auto/InvalidValue/2.49/bad/DHCPServerStaticLease_Address_255.256.255.255.target
[Unit]
Description=Foo
[DHCPServerStaticLease]
Address=255.256.255.255

