# tests/auto/InvalidValue/2.49/bad/DHCPServerStaticLease_Address_abc.target
[Unit]
Description=Foo
[DHCPServerStaticLease]
Address=abc

