# tests/auto/InvalidValue/2.49/bad/DHCPServerStaticLease_Address_0.0.0..target
[Unit]
Description=Foo
[DHCPServerStaticLease]
Address=0.0.0.

