# tests/auto/InvalidValue/2.49/bad/DHCPServer_RelayTarget_256.255.255.255.target
[Unit]
Description=Foo
[DHCPServer]
RelayTarget=256.255.255.255

