# tests/auto/InvalidValue/2.49/bad/DHCPServer_ServerAddress_255.256.255.255.target
[Unit]
Description=Foo
[DHCPServer]
ServerAddress=255.256.255.255

