# tests/auto/InvalidValue/2.49/bad/DHCPServer_RelayTarget_abc.target
[Unit]
Description=Foo
[DHCPServer]
RelayTarget=abc

