# tests/auto/InvalidValue/2.49/bad/DHCPServerStaticLease_Address_255.255.255.256.target
[Unit]
Description=Foo
[DHCPServerStaticLease]
Address=255.255.255.256

