# tests/auto/InvalidValue/2.49/bad/Network_ManageForeignRoutingPolicyRules_baz.target
[Unit]
Description=Foo
[Network]
ManageForeignRoutingPolicyRules=baz

