# tests/auto/InvalidValue/2.49/bad/Network_ManageForeignRoutingPolicyRules_bar.target
[Unit]
Description=Foo
[Network]
ManageForeignRoutingPolicyRules=bar

