# tests/auto/InvalidValue/2.49/bad/Network_ManageForeignRoutingPolicyRules_foo.target
[Unit]
Description=Foo
[Network]
ManageForeignRoutingPolicyRules=foo

