# tests/auto/OptionDeprecated/2.41/good/Service_BlockIOWriteBandwidthMax_--dev--foo%200K.target
[Unit]
Description=Foo
[Service]
BlockIOWriteBandwidthMax=123

