# tests/auto/OptionDeprecated/2.41/good/Service_BlockIOReadBandwidthMax_--dev--foo%200K.target
[Unit]
Description=Foo
[Service]
BlockIOReadBandwidthMax=123

