# tests/auto/OptionDeprecated/2.41/good/Service_BlockIOWeight_1.target
[Unit]
Description=Foo
[Service]
BlockIOWeight=123

