# tests/auto/OptionDeprecated/2.41/good/Socket_BlockIOReadBandwidthMax_--dev--foo%200K.target
[Unit]
Description=Foo
[Socket]
BlockIOReadBandwidthMax=123

