# tests/auto/OptionDeprecated/2.41/good/Socket_StartupCPUShares_2.target
[Unit]
Description=Foo
[Socket]
StartupCPUShares=123

