# tests/auto/OptionDeprecated/2.41/good/Swap_BlockIODeviceWeight_--dev--foo%201.target
[Unit]
Description=Foo
[Swap]
BlockIODeviceWeight=123

