# systemdlint/tests/auto/OptionTooNew/2.40/good/Socket_LogRateLimitIntervalSec_0sec.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=123

