# systemdlint/tests/auto/OptionTooNew/2.30/bad/Service_ReadWritePaths_abc.target
[Unit]
Description=Foo
[Service]
ReadWritePaths=123

