# systemdlint/tests/auto/OptionTooNew/2.30/bad/Service_ReadOnlyPaths_abc.target
[Unit]
Description=Foo
[Service]
ReadOnlyPaths=123

