# systemdlint/tests/auto/OptionTooNew/2.30/bad/Socket_InaccessableDirectories_abc.target
[Unit]
Description=Foo
[Socket]
InaccessableDirectories=123

