# systemdlint/tests/auto/OptionTooNew/2.30/bad/Swap_ReadOnlyPaths_abc.target
[Unit]
Description=Foo
[Swap]
ReadOnlyPaths=123

