# systemdlint/tests/auto/OptionTooNew/2.30/bad/Swap_MemoryDenyWriteExecute_yes.target
[Unit]
Description=Foo
[Swap]
MemoryDenyWriteExecute=123

