# systemdlint/tests/auto/OptionTooNew/2.30/bad/Swap_InaccessiblePaths_abc.target
[Unit]
Description=Foo
[Swap]
InaccessiblePaths=123

