# tests/manual/Multiplicity/2.31/bad/Service_OneShot_NoExecStart.service
[Unit]
Description=no ExecStart

[Service]
Type=oneshot