# systemdlint/tests/manual/NoFailureCheck/2.31/bad/Service_ExecStart_ls@.service
[Unit]
Description=Foo

[Service]
Type=oneshot
ExecStart=-/bin/ls %i
