# systemdlint/tests/auto/OptionTooNew/2.42/bad/BridgeFDB_AssociatedWith_use.target
[Unit]
Description=Foo
[BridgeFDB]
MACAddress=01:23:45:67:89:ab
AssociatedWith=123

