# systemdlint/tests/auto/OptionTooNew/2.42/bad/BridgeFDB_AssociatedWith_use.target
[Unit]
Description=Foo
[BridgeFDB]
MACAddress=abc
AssociatedWith=123

