# systemdlint/tests/auto/OptionTooNew/2.42/bad/Bridge_ProxyARPWiFi_yes.target
[Unit]
Description=Foo
[Bridge]
ProxyARPWiFi=123

