# systemdlint/tests/auto/OptionTooNew/2.42/bad/DHCPv4_KeepConfiguration_abc.target
[Unit]
Description=Foo
[DHCPv4]
KeepConfiguration=123

