# systemdlint/tests/auto/OptionTooNew/2.42/bad/DHCPv4_UseDomains_yes.target
[Unit]
Description=Foo
[DHCPv4]
UseDomains=123

