# systemdlint/tests/auto/OptionTooNew/2.42/bad/DHCPv6_DUIDRawData_abc.target
[Unit]
Description=Foo
[DHCPv6]
DUIDRawData=123

