# systemdlint/tests/auto/OptionTooNew/2.42/bad/FooOverUDP_Local_1.1.1.1--0.target
[Unit]
Description=Foo
[FooOverUDP]
Local=123

