# systemdlint/tests/auto/OptionTooNew/2.42/bad/GENEVE_IPDoNotFragment_yes.target
[Unit]
Description=Foo
[GENEVE]
IPDoNotFragment=123

