# systemdlint/tests/auto/OptionTooNew/2.42/bad/Neighbor_LinkLayerAddress_abc.target
[Unit]
Description=Foo
[Neighbor]
LinkLayerAddress=123

