# systemdlint/tests/auto/OptionTooNew/2.42/bad/Service_TimeoutAbortSec_0sec.target
[Unit]
Description=Foo
[Service]
TimeoutAbortSec=123

