# systemdlint/tests/auto/OptionTooNew/2.42/bad/Socket_DefaultMemoryLow_0K.target
[Unit]
Description=Foo
[Socket]
DefaultMemoryLow=123

