# systemdlint/tests/auto/OptionTooNew/2.42/bad/Swap_NUMAPolicy_default.target
[Unit]
Description=Foo
[Swap]
NUMAPolicy=123

