# systemdlint/tests/auto/OptionTooNew/2.42/bad/Unit_IPIngressFilterPath_abc.target
[Unit]
Description=Foo
Description=abc
IPIngressFilterPath=123

