# tests/auto/SettingRestricted/2.43/good/DHCPv4_UseRoutes_yesAnonymize%3D1.target
[Unit]
Description=Foo
[DHCPv4]
UseRoutes=yes

