# tests/auto/OptionDeprecated/2.13/bad/Socket_ControlGroupAttribute_abc.target
[Unit]
Description=Foo
[Socket]
ControlGroupAttribute=123

