# tests/auto/OptionDeprecated/2.13/bad/Swap_ControlGroupModify_abc.target
[Unit]
Description=Foo
[Swap]
ControlGroupModify=123

