# tests/auto/OptionDeprecated/2.13/bad/Swap_ControlGroupPersistent_abc.target
[Unit]
Description=Foo
[Swap]
ControlGroupPersistent=123

