# systemdlint/tests/auto/OptionTooNew/2.32/bad/Service_BindPaths_abc.target
[Unit]
Description=Foo
[Service]
BindPaths=123

