# tests/auto/OptionTooNew/2.30/bad/Socket_MemoryDenyWriteExecute_yes.target
[Unit]
Description=Foo
[Socket]
MemoryDenyWriteExecute=123

