# tests/auto/OptionTooNew/2.30/bad/Socket_ReadWritePaths_abc.target
[Unit]
Description=Foo
[Socket]
ReadWritePaths=123

