# tests/auto/OptionTooNew/2.30/bad/Socket_ReadWriteDirectories_abc.target
[Unit]
Description=Foo
[Socket]
ReadWriteDirectories=123

