# systemdlint/tests/auto/OptionDeprecated/2.13/bad/Mount_ControlGroupAttribute_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
ControlGroupAttribute=123

