# systemdlint/tests/auto/OptionDeprecated/2.13/bad/Service_ControlGroupModify_abc.target
[Unit]
Description=Foo
[Service]
ControlGroupModify=123

