# systemdlint/tests/auto/OptionDeprecated/2.13/bad/Socket_ControlGroupModify_abc.target
[Unit]
Description=Foo
[Socket]
ControlGroupModify=123

