# systemdlint/tests/auto/OptionDeprecated/2.13/bad/Swap_ControlGroupAttribute_abc.target
[Unit]
Description=Foo
[Swap]
ControlGroupAttribute=123

