# systemdlint/tests/auto/InvalidValue/2.39/bad/Manager_NoNewPrivileges_baz.target
[Unit]
Description=Foo
[Manager]
NoNewPrivileges=baz

