# systemdlint/tests/auto/InvalidValue/2.39/bad/Mount_PrivateMounts_foo.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
PrivateMounts=foo

