# systemdlint/tests/auto/InvalidValue/2.39/bad/Service_NoNewPrivileges_bar.target
[Unit]
Description=Foo
[Service]
NoNewPrivileges=bar

