# systemdlint/tests/auto/InvalidValue/2.39/bad/Service_PrivateMounts_foo.target
[Unit]
Description=Foo
[Service]
PrivateMounts=foo

