# systemdlint/tests/auto/InvalidValue/2.39/bad/Service_NoNewPrivileges_baz.target
[Unit]
Description=Foo
[Service]
NoNewPrivileges=baz

