# systemdlint/tests/auto/InvalidValue/2.39/bad/Socket_NoNewPrivileges_bar.target
[Unit]
Description=Foo
[Socket]
NoNewPrivileges=bar

