# systemdlint/tests/auto/InvalidValue/2.39/bad/Socket_NoNewPrivileges_baz.target
[Unit]
Description=Foo
[Socket]
NoNewPrivileges=baz

