# systemdlint/tests/auto/InvalidValue/2.39/bad/Socket_PrivateMounts_baz.target
[Unit]
Description=Foo
[Socket]
PrivateMounts=baz

