# systemdlint/tests/auto/InvalidValue/2.39/bad/Socket_NoNewPrivileges_foo.target
[Unit]
Description=Foo
[Socket]
NoNewPrivileges=foo

