# systemdlint/tests/auto/InvalidValue/2.39/bad/Swap_NoNewPrivileges_baz.target
[Unit]
Description=Foo
[Swap]
NoNewPrivileges=baz

