# systemdlint/tests/auto/InvalidValue/2.39/bad/Swap_PrivateMounts_foo.target
[Unit]
Description=Foo
[Swap]
PrivateMounts=foo

