# systemdlint/tests/auto/InvalidValue/2.33/good/Service_BindReadOnlyPaths_abc.target
[Unit]
Description=Foo
[Service]
BindReadOnlyPaths=abc

