# systemdlint/tests/auto/InvalidValue/2.33/good/Service_BindPaths_abc.target
[Unit]
Description=Foo
[Service]
BindPaths=abc

