# systemdlint/tests/auto/InvalidValue/2.33/good/Socket_BindReadOnlyPaths_abc.target
[Unit]
Description=Foo
[Socket]
BindReadOnlyPaths=abc

