# systemdlint/tests/auto/InvalidValue/2.33/good/Swap_BindPaths_abc.target
[Unit]
Description=Foo
[Swap]
BindPaths=abc

