# tests/auto/OptionDeprecated/2.30/good/Socket_Capabilities_abc.target
[Unit]
Description=Foo
[Socket]
Capabilities=123

