# tests/auto/OptionTooNew/2.14/good/Socket_ProtectSystem_full.target
[Unit]
Description=Foo
[Socket]
ProtectSystem=123

