# tests/auto/OptionTooNew/2.14/good/Socket_ProtectHome_tmpfs.target
[Unit]
Description=Foo
[Socket]
ProtectHome=123

