# tests/auto/OptionTooNew/2.14/good/Socket_SocketGroup_iamagroup.target
[Unit]
Description=Foo
[Socket]
SocketGroup=123

