# tests/auto/OptionTooNew/2.42/bad/DHCPv4_DUIDType_abc.target
[Unit]
Description=Foo
[DHCPv4]
DUIDType=123

