# tests/auto/OptionTooNew/2.42/bad/DHCP_KeepConfiguration_abc.target
[Unit]
Description=Foo
[DHCP]
KeepConfiguration=123

