# tests/auto/OptionTooNew/2.42/bad/DHCPv4_UseTimezone_yes.target
[Unit]
Description=Foo
[DHCPv4]
UseTimezone=123

