# tests/auto/OptionTooNew/2.42/bad/DHCPv6_UseTimezone_yes.target
[Unit]
Description=Foo
[DHCPv6]
UseTimezone=123

