# tests/auto/OptionTooNew/2.42/bad/DHCPv6_ClientIdentifier_mac.target
[Unit]
Description=Foo
[DHCPv6]
ClientIdentifier=123

