# tests/auto/OptionTooNew/2.42/bad/DHCPv6_KeepConfiguration_abc.target
[Unit]
Description=Foo
[DHCPv6]
KeepConfiguration=123

