# tests/auto/OptionTooNew/2.42/bad/DHCPv4_RequestBroadcast_abc.target
[Unit]
Description=Foo
[DHCPv4]
RequestBroadcast=123

