# tests/auto/OptionTooNew/2.42/bad/GENEVE_GenericProtocolExtension_yes.target
[Unit]
Description=Foo
[GENEVE]
GenericProtocolExtension=123

