# tests/auto/OptionTooNew/2.42/bad/Link_OriginalName_abc.target
[Unit]
Description=Foo
[Link]
OriginalName=123

