# tests/auto/OptionTooNew/2.42/bad/Manager_NUMAPolicy_default.target
[Unit]
Description=Foo
[Manager]
NUMAPolicy=123

