# tests/auto/OptionTooNew/2.42/bad/Match_Property_abc.target
[Unit]
Description=Foo
[Match]
Property=123

