# tests/auto/OptionTooNew/2.42/bad/Route_TTLPropagate_yes.target
[Unit]
Description=Foo
[Route]
TTLPropagate=123

