# tests/auto/OptionTooNew/2.42/bad/Route_FastOpenNoCookie_yes.target
[Unit]
Description=Foo
[Route]
FastOpenNoCookie=123

