# tests/auto/OptionTooNew/2.42/bad/Service_ExecCondition_abc.target
[Unit]
Description=Foo
[Service]
ExecCondition=123

