# tests/auto/OptionTooNew/2.42/bad/VXLAN_IPDoNotFragment_yes.target
[Unit]
Description=Foo
[VXLAN]
IPDoNotFragment=123

