# systemdlint/tests/auto/OptionTooNew/2.42/good/Service_NetworkNamespacePath_--tmp.target
[Unit]
Description=Foo
[Service]
NetworkNamespacePath=123

