# tests/auto/OptionDeprecated/2.31/bad/Service_Capabilities_abc.target
[Unit]
Description=Foo
[Service]
Capabilities=123

