# tests/auto/InvalidValue/2.31/good/Address_Scope_0.target
[Unit]
Description=Foo
[Address]
Address=abc
Scope=0

