# tests/auto/InvalidValue/2.31/good/Address_Peer_255.255.255.255--32.target
[Unit]
Description=Foo
[Address]
Address=abc
Peer=255.255.255.255/32

