# tests/auto/InvalidValue/2.31/good/Address_ManageTemporaryAddress_off.target
[Unit]
Description=Foo
[Address]
Address=abc
ManageTemporaryAddress=off

