# tests/auto/InvalidValue/2.31/good/Address_Scope_global.target
[Unit]
Description=Foo
[Address]
Address=abc
Scope=global

