# tests/auto/InvalidValue/2.31/good/Address_DuplicateAddressDetection_1.target
[Unit]
Description=Foo
[Address]
Address=abc
DuplicateAddressDetection=1

