# tests/auto/InvalidValue/2.31/good/Address_DuplicateAddressDetection_off.target
[Unit]
Description=Foo
[Address]
Address=abc
DuplicateAddressDetection=off

