# tests/auto/InvalidValue/2.31/good/Address_ManageTemporaryAddress_0.target
[Unit]
Description=Foo
[Address]
Address=abc
ManageTemporaryAddress=0

