# tests/auto/InvalidValue/2.31/good/Address_AutoJoin_1.target
[Unit]
Description=Foo
[Address]
Address=abc
AutoJoin=1

