# tests/auto/InvalidValue/2.31/good/Address_Address_abc.target
[Unit]
Description=Foo
[Address]
Address=abc
Address=abc

