# tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_21us.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=21us

