# tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_infinity.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=infinity

