# tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_5hr.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=5hr

