# tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_7minutes.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=7minutes

