# tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_19s.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=19s

