# tests/auto/InvalidValue/2.31/good/Automount_What_abc.target
[Unit]
Description=Foo
[Automount]
What=abc

