# tests/auto/InvalidValue/2.31/good/Automount_Where_--dev--.target
[Unit]
Description=Foo
[Automount]
Where=/dev/

