# tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_24w.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=24w

