# tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_13months.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=13months

