# tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_12month.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=12month

