# tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_4hour.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=4hour

