# tests/auto/InvalidValue/2.31/good/Automount_TimeoutIdleSec_3hours.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=3hours

