# tests/auto/InvalidValue/2.31/good/Bond_AdSelect_bandwidth.target
[Unit]
Description=Foo
[Bond]
AdSelect=bandwidth

