# tests/auto/InvalidValue/2.31/good/Bond_Mode_broadcast.target
[Unit]
Description=Foo
[Bond]
Mode=broadcast

