# tests/auto/InvalidValue/2.31/good/Bond_FailOverMACPolicy_follow.target
[Unit]
Description=Foo
[Bond]
FailOverMACPolicy=follow

