# tests/auto/InvalidValue/2.31/good/Bond_PrimaryReselectPolicy_better.target
[Unit]
Description=Foo
[Bond]
PrimaryReselectPolicy=better

