# tests/auto/InvalidValue/2.31/good/Bond_DynamicTransmitLoadBalancing_1.target
[Unit]
Description=Foo
[Bond]
DynamicTransmitLoadBalancing=1

