# tests/auto/InvalidValue/2.31/good/Bond_AdUserPortKey_1023.target
[Unit]
Description=Foo
[Bond]
AdUserPortKey=1023

