# tests/auto/InvalidValue/2.31/good/Bond_LACPTransmitRate_fast.target
[Unit]
Description=Foo
[Bond]
LACPTransmitRate=fast

