# tests/auto/InvalidValue/2.31/good/Bond_AllSlavesActive_no.target
[Unit]
Description=Foo
[Bond]
AllSlavesActive=no

