# tests/auto/InvalidValue/2.31/good/Bond_TransmitHashPolicy_encap2%2B3.target
[Unit]
Description=Foo
[Bond]
TransmitHashPolicy=encap2+3

