# tests/auto/InvalidValue/2.31/good/Bond_AdActorSystemPriority_1.target
[Unit]
Description=Foo
[Bond]
AdActorSystemPriority=1

