# tests/auto/InvalidValue/2.31/good/Bond_AllSlavesActive_true.target
[Unit]
Description=Foo
[Bond]
AllSlavesActive=true

