# tests/auto/InvalidValue/2.31/good/Bond_TransmitHashPolicy_layer2.target
[Unit]
Description=Foo
[Bond]
TransmitHashPolicy=layer2

