# tests/auto/InvalidValue/2.31/good/Bond_ARPIntervalSec_0.target
[Unit]
Description=Foo
[Bond]
ARPIntervalSec=0

