# tests/auto/InvalidValue/2.31/good/Bond_LACPTransmitRate_slow.target
[Unit]
Description=Foo
[Bond]
LACPTransmitRate=slow

