# tests/auto/InvalidValue/2.31/good/Bond_LearnPacketIntervalSec_2147483647.target
[Unit]
Description=Foo
[Bond]
LearnPacketIntervalSec=2147483647

