# tests/auto/InvalidValue/2.31/good/Bond_ARPIntervalSec_9999999999999.target
[Unit]
Description=Foo
[Bond]
ARPIntervalSec=9999999999999

