# tests/auto/InvalidValue/2.31/good/Bond_GratuitousARP_255.target
[Unit]
Description=Foo
[Bond]
GratuitousARP=255

