# tests/auto/InvalidValue/2.31/good/Bond_GratuitousARP_0.target
[Unit]
Description=Foo
[Bond]
GratuitousARP=0

