# tests/auto/InvalidValue/2.31/good/Bond_PacketsPerSlave_65535.target
[Unit]
Description=Foo
[Bond]
PacketsPerSlave=65535

