# tests/auto/InvalidValue/2.31/good/Bond_AllSlavesActive_off.target
[Unit]
Description=Foo
[Bond]
AllSlavesActive=off

