# tests/auto/InvalidValue/2.31/good/Bond_AllSlavesActive_on.target
[Unit]
Description=Foo
[Bond]
AllSlavesActive=on

