# tests/auto/InvalidValue/2.31/good/Bond_MinLinks_9999999999999.target
[Unit]
Description=Foo
[Bond]
MinLinks=9999999999999

