# tests/auto/InvalidValue/2.31/good/Bond_DynamicTransmitLoadBalancing_true.target
[Unit]
Description=Foo
[Bond]
DynamicTransmitLoadBalancing=true

