# tests/auto/InvalidValue/2.31/good/Bond_PrimaryReselectPolicy_failure.target
[Unit]
Description=Foo
[Bond]
PrimaryReselectPolicy=failure

