# tests/auto/InvalidValue/2.31/good/Bond_MIIMonitorSec_9999999999999.target
[Unit]
Description=Foo
[Bond]
MIIMonitorSec=9999999999999

