# tests/auto/InvalidValue/2.31/good/Bond_AdSelect_stable.target
[Unit]
Description=Foo
[Bond]
AdSelect=stable

