# tests/auto/InvalidValue/2.31/good/Bridge_MulticastToUnicast_off.target
[Unit]
Description=Foo
[Bridge]
MulticastToUnicast=off

