# tests/auto/InvalidValue/2.31/good/Bridge_UseBPDU_0.target
[Unit]
Description=Foo
[Bridge]
UseBPDU=0

