# tests/auto/InvalidValue/2.31/good/Bridge_UseBPDU_yes.target
[Unit]
Description=Foo
[Bridge]
UseBPDU=yes

