# tests/auto/InvalidValue/2.31/good/Bridge_Learning_no.target
[Unit]
Description=Foo
[Bridge]
Learning=no

