# tests/auto/InvalidValue/2.31/good/Bridge_UseBPDU_off.target
[Unit]
Description=Foo
[Bridge]
UseBPDU=off

