# tests/auto/InvalidValue/2.31/good/Bridge_UseBPDU_on.target
[Unit]
Description=Foo
[Bridge]
UseBPDU=on

