# tests/auto/InvalidValue/2.31/good/Bridge_MulticastFlood_false.target
[Unit]
Description=Foo
[Bridge]
MulticastFlood=false

