# tests/auto/InvalidValue/2.31/good/Bridge_UseBPDU_true.target
[Unit]
Description=Foo
[Bridge]
UseBPDU=true

