# tests/auto/InvalidValue/2.31/good/Bridge_MulticastFlood_no.target
[Unit]
Description=Foo
[Bridge]
MulticastFlood=no

