# tests/auto/InvalidValue/2.31/good/Bridge_NeighborSuppression_1.target
[Unit]
Description=Foo
[Bridge]
NeighborSuppression=1

