# tests/auto/InvalidValue/2.31/good/Bridge_HairPin_off.target
[Unit]
Description=Foo
[Bridge]
HairPin=off

