# tests/auto/InvalidValue/2.31/good/Bridge_MulticastToUnicast_true.target
[Unit]
Description=Foo
[Bridge]
MulticastToUnicast=true

