# tests/auto/InvalidValue/2.31/good/Bridge_MulticastToUnicast_false.target
[Unit]
Description=Foo
[Bridge]
MulticastToUnicast=false

