# tests/auto/InvalidValue/2.31/good/Bridge_HairPin_on.target
[Unit]
Description=Foo
[Bridge]
HairPin=on

