# tests/auto/InvalidValue/2.31/good/Bridge_MulticastToUnicast_0.target
[Unit]
Description=Foo
[Bridge]
MulticastToUnicast=0

