# tests/auto/InvalidValue/2.31/good/Bridge_Learning_true.target
[Unit]
Description=Foo
[Bridge]
Learning=true

