# tests/auto/InvalidValue/2.31/good/Bridge_Learning_1.target
[Unit]
Description=Foo
[Bridge]
Learning=1

