# tests/auto/InvalidValue/2.31/good/Bridge_NeighborSuppression_off.target
[Unit]
Description=Foo
[Bridge]
NeighborSuppression=off

