# tests/auto/InvalidValue/2.31/good/Bridge_HairPin_0.target
[Unit]
Description=Foo
[Bridge]
HairPin=0

