# tests/auto/InvalidValue/2.31/good/Bridge_AllowPortToBeRoot_yes.target
[Unit]
Description=Foo
[Bridge]
AllowPortToBeRoot=yes

