# tests/auto/InvalidValue/2.31/good/Bridge_Learning_false.target
[Unit]
Description=Foo
[Bridge]
Learning=false

