# tests/auto/InvalidValue/2.31/good/Bridge_HairPin_false.target
[Unit]
Description=Foo
[Bridge]
HairPin=false

