# tests/auto/InvalidValue/2.31/good/Bridge_UseBPDU_no.target
[Unit]
Description=Foo
[Bridge]
UseBPDU=no

