# tests/auto/InvalidValue/2.31/good/Bridge_Priority_0.target
[Unit]
Description=Foo
[Bridge]
Priority=0

