# tests/auto/InvalidValue/2.31/good/Bridge_Learning_on.target
[Unit]
Description=Foo
[Bridge]
Learning=on

