# tests/auto/InvalidValue/2.31/good/Bridge_UnicastFlood_yes.target
[Unit]
Description=Foo
[Bridge]
UnicastFlood=yes

