# tests/auto/InvalidValue/2.31/good/Bridge_NeighborSuppression_false.target
[Unit]
Description=Foo
[Bridge]
NeighborSuppression=false

