# tests/auto/InvalidValue/2.31/good/Bridge_FastLeave_false.target
[Unit]
Description=Foo
[Bridge]
FastLeave=false

