# tests/auto/InvalidValue/2.31/good/Bridge_FastLeave_off.target
[Unit]
Description=Foo
[Bridge]
FastLeave=off

