# tests/auto/InvalidValue/2.31/good/Bridge_Learning_off.target
[Unit]
Description=Foo
[Bridge]
Learning=off

