# tests/auto/InvalidValue/2.31/good/Bridge_UseBPDU_1.target
[Unit]
Description=Foo
[Bridge]
UseBPDU=1

