# tests/auto/InvalidValue/2.31/good/Bridge_UnicastFlood_off.target
[Unit]
Description=Foo
[Bridge]
UnicastFlood=off

