# tests/auto/InvalidValue/2.31/good/Bridge_AllowPortToBeRoot_0.target
[Unit]
Description=Foo
[Bridge]
AllowPortToBeRoot=0

