# tests/auto/InvalidValue/2.31/good/Bridge_NeighborSuppression_true.target
[Unit]
Description=Foo
[Bridge]
NeighborSuppression=true

