# tests/auto/InvalidValue/2.31/good/Bridge_NeighborSuppression_no.target
[Unit]
Description=Foo
[Bridge]
NeighborSuppression=no

