# tests/auto/InvalidValue/2.31/good/Bridge_Cost_65535.target
[Unit]
Description=Foo
[Bridge]
Cost=65535

