# tests/auto/InvalidValue/2.31/good/Bridge_FastLeave_no.target
[Unit]
Description=Foo
[Bridge]
FastLeave=no

