# tests/auto/InvalidValue/2.31/good/Bridge_UnicastFlood_0.target
[Unit]
Description=Foo
[Bridge]
UnicastFlood=0

