# tests/auto/InvalidValue/2.31/good/CAN_RestartSec_22weeks.target
[Unit]
Description=Foo
[CAN]
RestartSec=22weeks

