# tests/auto/InvalidValue/2.31/good/CAN_RestartSec_11m.target
[Unit]
Description=Foo
[CAN]
RestartSec=11m

