# tests/auto/InvalidValue/2.31/good/CAN_TripleSampling_on.target
[Unit]
Description=Foo
[CAN]
TripleSampling=on

