# tests/auto/InvalidValue/2.31/good/CAN_SamplePoint_abc.target
[Unit]
Description=Foo
[CAN]
SamplePoint=abc

