# tests/auto/InvalidValue/2.31/good/CAN_TripleSampling_yes.target
[Unit]
Description=Foo
[CAN]
TripleSampling=yes

