# tests/auto/InvalidValue/2.31/good/CAN_Bitrate_9999999999999K.target
[Unit]
Description=Foo
[CAN]
Bitrate=9999999999999K

