# tests/auto/InvalidValue/2.31/good/CAN_TripleSampling_true.target
[Unit]
Description=Foo
[CAN]
TripleSampling=true

