# tests/auto/InvalidValue/2.31/good/CAN_RestartSec_20usec.target
[Unit]
Description=Foo
[CAN]
RestartSec=20usec

