# tests/auto/InvalidValue/2.31/good/CAN_RestartSec_18sec.target
[Unit]
Description=Foo
[CAN]
RestartSec=18sec

