# tests/auto/InvalidValue/2.31/good/CAN_TripleSampling_false.target
[Unit]
Description=Foo
[CAN]
TripleSampling=false

