# tests/auto/InvalidValue/2.31/good/DHCP_Anonymize_false.target
[Unit]
Description=Foo
[DHCP]
Anonymize=false

