# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitRouter_on.target
[Unit]
Description=Foo
[DHCPServer]
EmitRouter=on

