# tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_6h.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=6h

