# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitDNS_1.target
[Unit]
Description=Foo
[DHCPServer]
EmitDNS=1

