# tests/auto/InvalidValue/2.31/good/DHCP_Anonymize_1.target
[Unit]
Description=Foo
[DHCP]
Anonymize=1

