# tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_16seconds.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=16seconds

