# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitTimezone_true.target
[Unit]
Description=Foo
[DHCPServer]
EmitTimezone=true

