# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitNTP_0.target
[Unit]
Description=Foo
[DHCPServer]
EmitNTP=0

