# tests/auto/InvalidValue/2.31/good/DHCP_UseDomains_route.target
[Unit]
Description=Foo
[DHCP]
UseDomains=route

