# tests/auto/InvalidValue/2.31/good/DHCP_CriticalConnection_true.target
[Unit]
Description=Foo
[DHCP]
CriticalConnection=true

