# tests/auto/InvalidValue/2.31/good/DHCP_UseMTU_0.target
[Unit]
Description=Foo
[DHCP]
UseMTU=0

