# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitDNS_true.target
[Unit]
Description=Foo
[DHCPServer]
EmitDNS=true

