# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitDNS_yes.target
[Unit]
Description=Foo
[DHCPServer]
EmitDNS=yes

