# tests/auto/InvalidValue/2.31/good/DHCP_RapidCommit_no.target
[Unit]
Description=Foo
[DHCP]
RapidCommit=no

