# tests/auto/InvalidValue/2.31/good/DHCP_SendHostname_false.target
[Unit]
Description=Foo
[DHCP]
SendHostname=false

