# tests/auto/InvalidValue/2.31/good/DHCP_UseRoutes_off.target
[Unit]
Description=Foo
[DHCP]
UseRoutes=off

