# tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_11m.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=11m

