# tests/auto/InvalidValue/2.31/good/DHCP_UseNTP_false.target
[Unit]
Description=Foo
[DHCP]
UseNTP=false

