# tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_infinity.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=infinity

