# tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_17second.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=17second

