# tests/auto/InvalidValue/2.31/good/DHCP_UseHostname_1.target
[Unit]
Description=Foo
[DHCP]
UseHostname=1

