# tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_14msec.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=14msec

