# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitNTP_1.target
[Unit]
Description=Foo
[DHCPServer]
EmitNTP=1

