# tests/auto/InvalidValue/2.31/good/DHCP_RapidCommit_1.target
[Unit]
Description=Foo
[DHCP]
RapidCommit=1

