# tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_infinity.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=infinity

