# tests/auto/InvalidValue/2.31/good/DHCP_UseTimezone_off.target
[Unit]
Description=Foo
[DHCP]
UseTimezone=off

