# tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_18sec.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=18sec

