# tests/auto/InvalidValue/2.31/good/DHCP_DUIDRawData_abc.target
[Unit]
Description=Foo
[DHCP]
DUIDRawData=abc

