# tests/auto/InvalidValue/2.31/good/DHCP_RapidCommit_off.target
[Unit]
Description=Foo
[DHCP]
RapidCommit=off

