# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitTimezone_off.target
[Unit]
Description=Foo
[DHCPServer]
EmitTimezone=off

