# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitDNS_off.target
[Unit]
Description=Foo
[DHCPServer]
EmitDNS=off

