# tests/auto/InvalidValue/2.31/good/DHCP_UseRoutes_true.target
[Unit]
Description=Foo
[DHCP]
UseRoutes=true

