# tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_3hours.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=3hours

