# tests/auto/InvalidValue/2.31/good/DHCP_DUIDType_abc.target
[Unit]
Description=Foo
[DHCP]
DUIDType=abc

