# tests/auto/InvalidValue/2.31/good/DHCP_UseDNS_0.target
[Unit]
Description=Foo
[DHCP]
UseDNS=0

