# tests/auto/InvalidValue/2.31/good/DHCP_CriticalConnection_1.target
[Unit]
Description=Foo
[DHCP]
CriticalConnection=1

