# tests/auto/InvalidValue/2.31/good/DHCPServer_NTP_abc.target
[Unit]
Description=Foo
[DHCPServer]
NTP=abc

