# tests/auto/InvalidValue/2.31/good/DHCPServer_PoolOffset_0.target
[Unit]
Description=Foo
[DHCPServer]
PoolOffset=0

