# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitNTP_yes.target
[Unit]
Description=Foo
[DHCPServer]
EmitNTP=yes

