# tests/auto/InvalidValue/2.31/good/DHCP_RapidCommit_yes.target
[Unit]
Description=Foo
[DHCP]
RapidCommit=yes

