# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitTimezone_0.target
[Unit]
Description=Foo
[DHCPServer]
EmitTimezone=0

