# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitTimezone_1.target
[Unit]
Description=Foo
[DHCPServer]
EmitTimezone=1

