# tests/auto/InvalidValue/2.31/good/DHCP_SendHostname_on.target
[Unit]
Description=Foo
[DHCP]
SendHostname=on

