# tests/auto/InvalidValue/2.31/good/DHCPServer_Timezone_abc.target
[Unit]
Description=Foo
[DHCPServer]
Timezone=abc

