# tests/auto/InvalidValue/2.31/good/DHCP_UseMTU_1.target
[Unit]
Description=Foo
[DHCP]
UseMTU=1

