# tests/auto/InvalidValue/2.31/good/DHCP_UserClass_abc.target
[Unit]
Description=Foo
[DHCP]
UserClass=abc

