# tests/auto/InvalidValue/2.31/good/DHCPServer_DefaultLeaseTimeSec_23week.target
[Unit]
Description=Foo
[DHCPServer]
DefaultLeaseTimeSec=23week

