# tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_4hour.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=4hour

