# tests/auto/InvalidValue/2.31/good/DHCP_RouteMetric_0.target
[Unit]
Description=Foo
[DHCP]
RouteMetric=0

