# tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_25years.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=25years

