# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitTimezone_yes.target
[Unit]
Description=Foo
[DHCPServer]
EmitTimezone=yes

