# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitRouter_0.target
[Unit]
Description=Foo
[DHCPServer]
EmitRouter=0

