# tests/auto/InvalidValue/2.31/good/DHCP_UseDNS_true.target
[Unit]
Description=Foo
[DHCP]
UseDNS=true

