# tests/auto/InvalidValue/2.31/good/DHCP_UseDomains_no.target
[Unit]
Description=Foo
[DHCP]
UseDomains=no

