# tests/auto/InvalidValue/2.31/good/DHCP_UseNTP_true.target
[Unit]
Description=Foo
[DHCP]
UseNTP=true

