# tests/auto/InvalidValue/2.31/good/DHCP_UseDNS_false.target
[Unit]
Description=Foo
[DHCP]
UseDNS=false

