# tests/auto/InvalidValue/2.31/good/DHCP_UseRoutes_1.target
[Unit]
Description=Foo
[DHCP]
UseRoutes=1

