# tests/auto/InvalidValue/2.31/good/DHCP_UseMTU_no.target
[Unit]
Description=Foo
[DHCP]
UseMTU=no

