# tests/auto/InvalidValue/2.31/good/DHCP_UseNTP_1.target
[Unit]
Description=Foo
[DHCP]
UseNTP=1

