# tests/auto/InvalidValue/2.31/good/DHCP_UseMTU_true.target
[Unit]
Description=Foo
[DHCP]
UseMTU=true

