# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitTimezone_false.target
[Unit]
Description=Foo
[DHCPServer]
EmitTimezone=false

