# tests/auto/InvalidValue/2.31/good/DHCP_SendHostname_yes.target
[Unit]
Description=Foo
[DHCP]
SendHostname=yes

