# tests/auto/InvalidValue/2.31/good/DHCPServer_MaxLeaseTimeSec_1day.target
[Unit]
Description=Foo
[DHCPServer]
MaxLeaseTimeSec=1day

