# tests/auto/InvalidValue/2.31/good/DHCP_UseDomains_1.target
[Unit]
Description=Foo
[DHCP]
UseDomains=1

