# tests/auto/InvalidValue/2.31/good/DHCP_UseNTP_no.target
[Unit]
Description=Foo
[DHCP]
UseNTP=no

