# tests/auto/InvalidValue/2.31/good/DHCP_UseNTP_off.target
[Unit]
Description=Foo
[DHCP]
UseNTP=off

