# tests/auto/InvalidValue/2.31/good/DHCP_Anonymize_on.target
[Unit]
Description=Foo
[DHCP]
Anonymize=on

