# tests/auto/InvalidValue/2.31/good/DHCP_UseRoutes_no.target
[Unit]
Description=Foo
[DHCP]
UseRoutes=no

