# tests/auto/InvalidValue/2.31/good/DHCP_UseTimezone_on.target
[Unit]
Description=Foo
[DHCP]
UseTimezone=on

