# tests/auto/InvalidValue/2.31/good/DHCP_UseMTU_on.target
[Unit]
Description=Foo
[DHCP]
UseMTU=on

