# tests/auto/InvalidValue/2.31/good/DHCP_ListenPort_65535.target
[Unit]
Description=Foo
[DHCP]
ListenPort=65535

