# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitDNS_0.target
[Unit]
Description=Foo
[DHCPServer]
EmitDNS=0

