# tests/auto/InvalidValue/2.31/good/DHCP_UseDNS_on.target
[Unit]
Description=Foo
[DHCP]
UseDNS=on

