# tests/auto/InvalidValue/2.31/good/DHCP_UseHostname_false.target
[Unit]
Description=Foo
[DHCP]
UseHostname=false

