# tests/auto/InvalidValue/2.31/good/DHCP_UseDNS_1.target
[Unit]
Description=Foo
[DHCP]
UseDNS=1

