# tests/auto/InvalidValue/2.31/good/DHCPServer_PoolOffset_9999999999999.target
[Unit]
Description=Foo
[DHCPServer]
PoolOffset=9999999999999

