# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitRouter_no.target
[Unit]
Description=Foo
[DHCPServer]
EmitRouter=no

