# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitRouter_false.target
[Unit]
Description=Foo
[DHCPServer]
EmitRouter=false

