# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitTimezone_no.target
[Unit]
Description=Foo
[DHCPServer]
EmitTimezone=no

