# tests/auto/InvalidValue/2.31/good/DHCP_UseDomains_on.target
[Unit]
Description=Foo
[DHCP]
UseDomains=on

