# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitNTP_no.target
[Unit]
Description=Foo
[DHCPServer]
EmitNTP=no

