# tests/auto/InvalidValue/2.31/good/DHCP_SendHostname_no.target
[Unit]
Description=Foo
[DHCP]
SendHostname=no

