# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitRouter_yes.target
[Unit]
Description=Foo
[DHCPServer]
EmitRouter=yes

