# tests/auto/InvalidValue/2.31/good/DHCP_RequestBroadcast_abc.target
[Unit]
Description=Foo
[DHCP]
RequestBroadcast=abc

