# tests/auto/InvalidValue/2.31/good/DHCP_UseDNS_off.target
[Unit]
Description=Foo
[DHCP]
UseDNS=off

