# tests/auto/InvalidValue/2.31/good/DHCPServer_EmitNTP_false.target
[Unit]
Description=Foo
[DHCPServer]
EmitNTP=false

