# tests/auto/InvalidValue/2.31/good/DHCP_CriticalConnection_off.target
[Unit]
Description=Foo
[DHCP]
CriticalConnection=off

