# tests/auto/InvalidValue/2.31/good/DHCP_UseNTP_yes.target
[Unit]
Description=Foo
[DHCP]
UseNTP=yes

