# tests/auto/InvalidValue/2.31/good/FooOverUDP_Protocol_ipip.target
[Unit]
Description=Foo
[FooOverUDP]
Protocol=ipip

