# tests/auto/InvalidValue/2.31/good/FooOverUDP_Protocol_1.target
[Unit]
Description=Foo
[FooOverUDP]
Protocol=1

