# tests/auto/InvalidValue/2.31/good/GENEVE_UDPChecksum_on.target
[Unit]
Description=Foo
[GENEVE]
UDPChecksum=on

