# tests/auto/InvalidValue/2.31/good/GENEVE_UDP6ZeroChecksumTx_true.target
[Unit]
Description=Foo
[GENEVE]
UDP6ZeroChecksumTx=true

