# tests/auto/InvalidValue/2.31/good/GENEVE_TOS_255.target
[Unit]
Description=Foo
[GENEVE]
TOS=255

