# tests/auto/InvalidValue/2.31/good/GENEVE_FlowLabel_abc.target
[Unit]
Description=Foo
[GENEVE]
FlowLabel=abc

