# tests/auto/InvalidValue/2.31/good/GENEVE_UDP6ZeroChecksumTx_0.target
[Unit]
Description=Foo
[GENEVE]
UDP6ZeroChecksumTx=0

