# tests/auto/InvalidValue/2.31/good/Install_DefaultInstance_abc.target
[Unit]
Description=Foo
[Install]
DefaultInstance=abc

