# tests/auto/InvalidValue/2.31/good/IPv6Prefix_PreferredLifetimeSec_18sec.target
[Unit]
Description=Foo
[IPv6Prefix]
PreferredLifetimeSec=18sec

