# tests/auto/InvalidValue/2.31/good/IPv6Prefix_PreferredLifetimeSec_14msec.target
[Unit]
Description=Foo
[IPv6Prefix]
PreferredLifetimeSec=14msec

