# tests/auto/InvalidValue/2.31/good/IPv6AcceptRA_RouteTable_4294967295.target
[Unit]
Description=Foo
[IPv6AcceptRA]
RouteTable=4294967295

