# tests/auto/InvalidValue/2.31/good/IPv6Prefix_OnLink_yes.target
[Unit]
Description=Foo
[IPv6Prefix]
OnLink=yes

