# tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_EmitDomains_off.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
EmitDomains=off

