# tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_Managed_no.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
Managed=no

