# tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_DNSLifetimeSec_19s.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
DNSLifetimeSec=19s

