# tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_Domains_abc.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
Domains=abc

