# tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_OtherInformation_false.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
OtherInformation=false

