# tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_Managed_on.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
Managed=on

