# tests/auto/InvalidValue/2.31/good/IPv6Prefix_PreferredLifetimeSec_23week.target
[Unit]
Description=Foo
[IPv6Prefix]
PreferredLifetimeSec=23week

