# tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_RouterLifetimeSec_20usec.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
RouterLifetimeSec=20usec

