# tests/auto/InvalidValue/2.31/good/IPv6AcceptRA_RouteTable_1.target
[Unit]
Description=Foo
[IPv6AcceptRA]
RouteTable=1

