# tests/auto/InvalidValue/2.31/good/IPv6AddressLabel_Prefix_abc.target
[Unit]
Description=Foo
[IPv6AddressLabel]
Label=abc
Prefix=abc
Prefix=abc

