# tests/auto/InvalidValue/2.31/good/IPv6Prefix_PreferredLifetimeSec_1day.target
[Unit]
Description=Foo
[IPv6Prefix]
PreferredLifetimeSec=1day

