# tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_RouterPreference_high.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
RouterPreference=high

