# tests/auto/InvalidValue/2.31/good/IPv6AcceptRA_UseAutonomousPrefix_no.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseAutonomousPrefix=no

