# tests/auto/InvalidValue/2.31/good/IPv6PrefixDelegation_DNSLifetimeSec_16seconds.target
[Unit]
Description=Foo
[IPv6PrefixDelegation]
DNSLifetimeSec=16seconds

