# tests/auto/InvalidValue/2.31/good/IPv6AcceptRA_UseOnLinkPrefix_0.target
[Unit]
Description=Foo
[IPv6AcceptRA]
UseOnLinkPrefix=0

