# tests/auto/InvalidValue/2.31/good/IPVLAN_Flags_private.target
[Unit]
Description=Foo
[IPVLAN]
Flags=private

