# tests/auto/InvalidValue/2.31/good/IPVLAN_Mode_L3S.target
[Unit]
Description=Foo
[IPVLAN]
Mode=L3S

