# tests/auto/InvalidValue/2.31/good/IPVLAN_Flags_vepa.target
[Unit]
Description=Foo
[IPVLAN]
Flags=vepa

