# tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_18sec.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=18sec

