# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelKMsg_debug.target
[Unit]
Description=Foo
[Journal]
MaxLevelKMsg=debug

