# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelWall_emerg.target
[Unit]
Description=Foo
[Journal]
MaxLevelWall=emerg

