# tests/auto/InvalidValue/2.31/good/Journal_LineMax_79T.target
[Unit]
Description=Foo
[Journal]
LineMax=79T

