# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelKMsg_alert.target
[Unit]
Description=Foo
[Journal]
MaxLevelKMsg=alert

