# tests/auto/InvalidValue/2.31/good/Journal_SystemMaxUse_9999999999999K.target
[Unit]
Description=Foo
[Journal]
SystemMaxUse=9999999999999K

