# tests/auto/InvalidValue/2.31/good/Journal_ForwardToKMsg_true.target
[Unit]
Description=Foo
[Journal]
ForwardToKMsg=true

