# tests/auto/InvalidValue/2.31/good/Journal_ForwardToConsole_no.target
[Unit]
Description=Foo
[Journal]
ForwardToConsole=no

