# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelWall_info.target
[Unit]
Description=Foo
[Journal]
MaxLevelWall=info

