# tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_22weeks.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=22weeks

