# tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_3hours.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=3hours

