# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelKMsg_notice.target
[Unit]
Description=Foo
[Journal]
MaxLevelKMsg=notice

