# tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_23week.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=23week

