# tests/auto/InvalidValue/2.31/good/Journal_ReadKMsg_1.target
[Unit]
Description=Foo
[Journal]
ReadKMsg=1

