# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelConsole_alert.target
[Unit]
Description=Foo
[Journal]
MaxLevelConsole=alert

