# tests/auto/InvalidValue/2.31/good/Journal_SystemMaxFileSize_9999999999999G.target
[Unit]
Description=Foo
[Journal]
SystemMaxFileSize=9999999999999G

