# tests/auto/InvalidValue/2.31/good/Journal_ForwardToConsole_1.target
[Unit]
Description=Foo
[Journal]
ForwardToConsole=1

