# tests/auto/InvalidValue/2.31/good/Journal_RuntimeMaxUse_9999999999999G.target
[Unit]
Description=Foo
[Journal]
RuntimeMaxUse=9999999999999G

