# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelConsole_debug.target
[Unit]
Description=Foo
[Journal]
MaxLevelConsole=debug

