# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelWall_0.target
[Unit]
Description=Foo
[Journal]
MaxLevelWall=0

