# tests/auto/InvalidValue/2.31/good/Journal_Seal_1.target
[Unit]
Description=Foo
[Journal]
Seal=1

