# tests/auto/InvalidValue/2.31/good/Journal_RuntimeMaxFileSize_0T.target
[Unit]
Description=Foo
[Journal]
RuntimeMaxFileSize=0T

