# tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_0days.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=0days

