# tests/auto/InvalidValue/2.31/good/Journal_SystemKeepFree_0G.target
[Unit]
Description=Foo
[Journal]
SystemKeepFree=0G

