# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelConsole_crit.target
[Unit]
Description=Foo
[Journal]
MaxLevelConsole=crit

