# tests/auto/InvalidValue/2.31/good/Journal_ForwardToSyslog_no.target
[Unit]
Description=Foo
[Journal]
ForwardToSyslog=no

