# tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_27y.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=27y

