# tests/auto/InvalidValue/2.31/good/Journal_Compress_1.target
[Unit]
Description=Foo
[Journal]
Compress=1

