# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelWall_alert.target
[Unit]
Description=Foo
[Journal]
MaxLevelWall=alert

