# tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_8minute.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=8minute

