# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelKMsg_crit.target
[Unit]
Description=Foo
[Journal]
MaxLevelKMsg=crit

