# tests/auto/InvalidValue/2.31/good/Journal_MaxFileSec_9999999999999h.target
[Unit]
Description=Foo
[Journal]
MaxFileSec=9999999999999h

