# tests/auto/InvalidValue/2.31/good/Journal_RateLimitIntervalSec_0us.target
[Unit]
Description=Foo
[Journal]
RateLimitIntervalSec=0us

