# tests/auto/InvalidValue/2.31/good/Journal_RateLimitIntervalSec_0s.target
[Unit]
Description=Foo
[Journal]
RateLimitIntervalSec=0s

