# tests/auto/InvalidValue/2.31/good/Journal_ForwardToWall_on.target
[Unit]
Description=Foo
[Journal]
ForwardToWall=on

