# tests/auto/InvalidValue/2.31/good/Journal_RuntimeMaxUse_0K.target
[Unit]
Description=Foo
[Journal]
RuntimeMaxUse=0K

