# tests/auto/InvalidValue/2.31/good/Journal_Compress_off.target
[Unit]
Description=Foo
[Journal]
Compress=off

