# tests/auto/InvalidValue/2.31/good/Journal_ForwardToSyslog_false.target
[Unit]
Description=Foo
[Journal]
ForwardToSyslog=false

