# tests/auto/InvalidValue/2.31/good/Journal_ForwardToConsole_on.target
[Unit]
Description=Foo
[Journal]
ForwardToConsole=on

