# tests/auto/InvalidValue/2.31/good/Journal_RuntimeMaxFiles_0T.target
[Unit]
Description=Foo
[Journal]
RuntimeMaxFiles=0T

