# tests/auto/InvalidValue/2.31/good/Journal_Seal_off.target
[Unit]
Description=Foo
[Journal]
Seal=off

