# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelStore_info.target
[Unit]
Description=Foo
[Journal]
MaxLevelStore=info

