# tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_28%C2%B5s.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=28µs

