# tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_21us.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=21us

