# tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_5hr.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=5hr

