# tests/auto/InvalidValue/2.31/good/Journal_ForwardToSyslog_on.target
[Unit]
Description=Foo
[Journal]
ForwardToSyslog=on

