# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelKMsg_warning.target
[Unit]
Description=Foo
[Journal]
MaxLevelKMsg=warning

