# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelConsole_info.target
[Unit]
Description=Foo
[Journal]
MaxLevelConsole=info

