# tests/auto/InvalidValue/2.31/good/Journal_SystemMaxFileSize_0K.target
[Unit]
Description=Foo
[Journal]
SystemMaxFileSize=0K

