# tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_6h.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=6h

