# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelStore_notice.target
[Unit]
Description=Foo
[Journal]
MaxLevelStore=notice

