# tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_4hour.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=4hour

