# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelConsole_0.target
[Unit]
Description=Foo
[Journal]
MaxLevelConsole=0

