# tests/auto/InvalidValue/2.31/good/Journal_Seal_true.target
[Unit]
Description=Foo
[Journal]
Seal=true

