# tests/auto/InvalidValue/2.31/good/Journal_TTYPath_--dev--foo.target
[Unit]
Description=Foo
[Journal]
TTYPath=/dev/foo

