# tests/auto/InvalidValue/2.31/good/Journal_ForwardToWall_off.target
[Unit]
Description=Foo
[Journal]
ForwardToWall=off

