# tests/auto/InvalidValue/2.31/good/Journal_Storage_persistent.target
[Unit]
Description=Foo
[Journal]
Storage=persistent

