# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelConsole_err.target
[Unit]
Description=Foo
[Journal]
MaxLevelConsole=err

