# tests/auto/InvalidValue/2.31/good/Journal_RuntimeKeepFree_0M.target
[Unit]
Description=Foo
[Journal]
RuntimeKeepFree=0M

