# tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_9min.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=9min

