# tests/auto/InvalidValue/2.31/good/Journal_RuntimeMaxFileSize_9999999999999M.target
[Unit]
Description=Foo
[Journal]
RuntimeMaxFileSize=9999999999999M

