# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelSyslog_crit.target
[Unit]
Description=Foo
[Journal]
MaxLevelSyslog=crit

