# tests/auto/InvalidValue/2.31/good/Journal_Seal_no.target
[Unit]
Description=Foo
[Journal]
Seal=no

