# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelStore_err.target
[Unit]
Description=Foo
[Journal]
MaxLevelStore=err

