# tests/auto/InvalidValue/2.31/good/Journal_ForwardToWall_yes.target
[Unit]
Description=Foo
[Journal]
ForwardToWall=yes

