# tests/auto/InvalidValue/2.31/good/Journal_Compress_true.target
[Unit]
Description=Foo
[Journal]
Compress=true

