# tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_25years.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=25years

