# tests/auto/InvalidValue/2.31/good/Journal_Seal_yes.target
[Unit]
Description=Foo
[Journal]
Seal=yes

