# tests/auto/InvalidValue/2.31/good/Journal_ForwardToConsole_yes.target
[Unit]
Description=Foo
[Journal]
ForwardToConsole=yes

