# tests/auto/InvalidValue/2.31/good/Journal_SplitMode_none.target
[Unit]
Description=Foo
[Journal]
SplitMode=none

