# tests/auto/InvalidValue/2.31/good/Journal_ReadKMsg_yes.target
[Unit]
Description=Foo
[Journal]
ReadKMsg=yes

