# tests/auto/InvalidValue/2.31/good/Journal_ForwardToSyslog_true.target
[Unit]
Description=Foo
[Journal]
ForwardToSyslog=true

