# tests/auto/InvalidValue/2.31/good/Journal_Compress_yes.target
[Unit]
Description=Foo
[Journal]
Compress=yes

