# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelStore_debug.target
[Unit]
Description=Foo
[Journal]
MaxLevelStore=debug

