# tests/auto/InvalidValue/2.31/good/Journal_MaxLevelConsole_notice.target
[Unit]
Description=Foo
[Journal]
MaxLevelConsole=notice

