# tests/auto/InvalidValue/2.31/good/Journal_ReadKMsg_true.target
[Unit]
Description=Foo
[Journal]
ReadKMsg=true

