# tests/auto/InvalidValue/2.31/good/Journal_SyncIntervalSec_13months.target
[Unit]
Description=Foo
[Journal]
SyncIntervalSec=13months

