# tests/auto/InvalidValue/2.31/good/L2TP_Local_1.1.1.1--0.target
[Unit]
Description=Foo
[L2TP]
Local=1.1.1.1/0

