# tests/auto/InvalidValue/2.31/good/L2TP_Local_auto.target
[Unit]
Description=Foo
[L2TP]
Local=auto

