# tests/auto/InvalidValue/2.31/good/L2TP_Local_255.255.255.255.target
[Unit]
Description=Foo
[L2TP]
Local=255.255.255.255

