# tests/auto/InvalidValue/2.31/good/L2TP_Remote_255.255.255.255.target
[Unit]
Description=Foo
[L2TP]
Remote=255.255.255.255

