# tests/auto/InvalidValue/2.31/good/L2TP_EncapsulationType_ip.target
[Unit]
Description=Foo
[L2TP]
EncapsulationType=ip

