# tests/auto/InvalidValue/2.31/good/Link_AllMulticast_true.target
[Unit]
Description=Foo
[Link]
AllMulticast=true

