# tests/auto/InvalidValue/2.31/good/Link_Multicast_on.target
[Unit]
Description=Foo
[Link]
Multicast=on

