# tests/auto/InvalidValue/2.31/good/Link_RequiredForOnline_0.target
[Unit]
Description=Foo
[Link]
RequiredForOnline=0

