# tests/auto/InvalidValue/2.31/good/Link_RequiredForOnline_false.target
[Unit]
Description=Foo
[Link]
RequiredForOnline=false

