# tests/auto/InvalidValue/2.31/good/Link_RequiredForOnline_1.target
[Unit]
Description=Foo
[Link]
RequiredForOnline=1

