# tests/auto/InvalidValue/2.31/good/Link_Unmanaged_true.target
[Unit]
Description=Foo
[Link]
Unmanaged=true

