# tests/auto/InvalidValue/2.31/good/Link_MTUBytes_0G.target
[Unit]
Description=Foo
[Link]
MTUBytes=0G

