# tests/auto/InvalidValue/2.31/good/Link_ARP_true.target
[Unit]
Description=Foo
[Link]
ARP=true

