# tests/auto/InvalidValue/2.31/good/Link_RequiredForOnline_on.target
[Unit]
Description=Foo
[Link]
RequiredForOnline=on

