# tests/auto/InvalidValue/2.31/good/Link_Unmanaged_yes.target
[Unit]
Description=Foo
[Link]
Unmanaged=yes

