# tests/auto/InvalidValue/2.31/good/Link_MTUBytes_0K.target
[Unit]
Description=Foo
[Link]
MTUBytes=0K

