# tests/auto/InvalidValue/2.31/good/Link_Unmanaged_off.target
[Unit]
Description=Foo
[Link]
Unmanaged=off

