# tests/auto/InvalidValue/2.31/good/Link_RequiredForOnline_yes.target
[Unit]
Description=Foo
[Link]
RequiredForOnline=yes

