# tests/auto/InvalidValue/2.31/good/Link_Multicast_0.target
[Unit]
Description=Foo
[Link]
Multicast=0

