# tests/auto/InvalidValue/2.31/good/Login_HandleHibernateKey_halt.target
[Unit]
Description=Foo
[Login]
HandleHibernateKey=halt

