# tests/auto/InvalidValue/2.31/good/Login_SuspendKeyIgnoreInhibited_no.target
[Unit]
Description=Foo
[Login]
SuspendKeyIgnoreInhibited=no

