# tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_28%C2%B5s.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=28µs

