# tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_8minute.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=8minute

