# tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_18sec.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=18sec

