# tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_2d.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=2d

