# tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_16seconds.target
[Unit]
Description=Foo
[Login]
IdleActionSec=16seconds

