# tests/auto/InvalidValue/2.31/good/Login_LidSwitchIgnoreInhibited_off.target
[Unit]
Description=Foo
[Login]
LidSwitchIgnoreInhibited=off

