# tests/auto/InvalidValue/2.31/good/Login_HandleSuspendKey_suspend.target
[Unit]
Description=Foo
[Login]
HandleSuspendKey=suspend

