# tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_27y.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=27y

