# tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_7minutes.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=7minutes

