# tests/auto/InvalidValue/2.31/good/Login_InhibitorsMax_9999999999999K.target
[Unit]
Description=Foo
[Login]
InhibitorsMax=9999999999999K

