# tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_7minutes.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=7minutes

