# tests/auto/InvalidValue/2.31/good/Login_NAutoVTs_0.target
[Unit]
Description=Foo
[Login]
NAutoVTs=0

