# tests/auto/InvalidValue/2.31/good/Login_SuspendKeyIgnoreInhibited_true.target
[Unit]
Description=Foo
[Login]
SuspendKeyIgnoreInhibited=true

