# tests/auto/InvalidValue/2.31/good/Login_IdleAction_reboot.target
[Unit]
Description=Foo
[Login]
IdleAction=reboot

