# tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_10M.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=10M

