# tests/auto/InvalidValue/2.31/good/Login_SessionsMax_0K.target
[Unit]
Description=Foo
[Login]
SessionsMax=0K

