# tests/auto/InvalidValue/2.31/good/Login_NAutoVTs_9999999999999.target
[Unit]
Description=Foo
[Login]
NAutoVTs=9999999999999

