# tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_6h.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=6h

