# tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitch_lock.target
[Unit]
Description=Foo
[Login]
HandleLidSwitch=lock

