# tests/auto/InvalidValue/2.31/good/Login_IdleAction_ignore.target
[Unit]
Description=Foo
[Login]
IdleAction=ignore

