# tests/auto/InvalidValue/2.31/good/Login_ReserveVT_9999999999999.target
[Unit]
Description=Foo
[Login]
ReserveVT=9999999999999

