# tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitchExternalPower_suspend.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchExternalPower=suspend

