# tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_9min.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=9min

