# tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_16seconds.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=16seconds

