# tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_3hours.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=3hours

