# tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_24w.target
[Unit]
Description=Foo
[Login]
IdleActionSec=24w

