# tests/auto/InvalidValue/2.31/good/Login_HandleSuspendKey_hybrid-sleep.target
[Unit]
Description=Foo
[Login]
HandleSuspendKey=hybrid-sleep

