# tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_23week.target
[Unit]
Description=Foo
[Login]
IdleActionSec=23week

