# tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_infinity.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=infinity

