# tests/auto/InvalidValue/2.31/good/Login_KillUserProcesses_off.target
[Unit]
Description=Foo
[Login]
KillUserProcesses=off

