# tests/auto/InvalidValue/2.31/good/Login_HandleHibernateKey_reboot.target
[Unit]
Description=Foo
[Login]
HandleHibernateKey=reboot

