# tests/auto/InvalidValue/2.31/good/Login_RemoveIPC_0.target
[Unit]
Description=Foo
[Login]
RemoveIPC=0

