# tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitchDocked_lock.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchDocked=lock

