# tests/auto/InvalidValue/2.31/good/Login_HandleHibernateKey_kexec.target
[Unit]
Description=Foo
[Login]
HandleHibernateKey=kexec

