# tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_2d.target
[Unit]
Description=Foo
[Login]
IdleActionSec=2d

