# tests/auto/InvalidValue/2.31/good/Login_HibernateKeyIgnoreInhibited_on.target
[Unit]
Description=Foo
[Login]
HibernateKeyIgnoreInhibited=on

