# tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitchExternalPower_kexec.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchExternalPower=kexec

