# tests/auto/InvalidValue/2.31/good/Login_KillUserProcesses_on.target
[Unit]
Description=Foo
[Login]
KillUserProcesses=on

