# tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_23week.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=23week

