# tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_27y.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=27y

