# tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitchExternalPower_halt.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchExternalPower=halt

