# tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_1day.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=1day

