# tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_25years.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=25years

