# tests/auto/InvalidValue/2.31/good/Login_KillUserProcesses_true.target
[Unit]
Description=Foo
[Login]
KillUserProcesses=true

