# tests/auto/InvalidValue/2.31/good/Login_KillOnlyUsers_abc.target
[Unit]
Description=Foo
[Login]
KillOnlyUsers=abc

