# tests/auto/InvalidValue/2.31/good/Login_IdleAction_hibernate.target
[Unit]
Description=Foo
[Login]
IdleAction=hibernate

