# tests/auto/InvalidValue/2.31/good/Login_InhibitorsMax_0K.target
[Unit]
Description=Foo
[Login]
InhibitorsMax=0K

