# tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_17second.target
[Unit]
Description=Foo
[Login]
IdleActionSec=17second

