# tests/auto/InvalidValue/2.31/good/Login_HandlePowerKey_halt.target
[Unit]
Description=Foo
[Login]
HandlePowerKey=halt

