# tests/auto/InvalidValue/2.31/good/Login_LidSwitchIgnoreInhibited_yes.target
[Unit]
Description=Foo
[Login]
LidSwitchIgnoreInhibited=yes

