# tests/auto/InvalidValue/2.31/good/Login_LidSwitchIgnoreInhibited_on.target
[Unit]
Description=Foo
[Login]
LidSwitchIgnoreInhibited=on

