# tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_22weeks.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=22weeks

