# tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_11m.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=11m

