# tests/auto/InvalidValue/2.31/good/Login_KillUserProcesses_0.target
[Unit]
Description=Foo
[Login]
KillUserProcesses=0

