# tests/auto/InvalidValue/2.31/good/Login_LidSwitchIgnoreInhibited_no.target
[Unit]
Description=Foo
[Login]
LidSwitchIgnoreInhibited=no

