# tests/auto/InvalidValue/2.31/good/Login_PowerKeyIgnoreInhibited_yes.target
[Unit]
Description=Foo
[Login]
PowerKeyIgnoreInhibited=yes

