# tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_9min.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=9min

