# tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitchExternalPower_poweroff.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchExternalPower=poweroff

