# tests/auto/InvalidValue/2.31/good/Login_SessionsMax_9999999999999T.target
[Unit]
Description=Foo
[Login]
SessionsMax=9999999999999T

