# tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_7minutes.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=7minutes

