# tests/auto/InvalidValue/2.31/good/Login_IdleAction_hybrid-sleep.target
[Unit]
Description=Foo
[Login]
IdleAction=hybrid-sleep

