# tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_16seconds.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=16seconds

