# tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitch_poweroff.target
[Unit]
Description=Foo
[Login]
HandleLidSwitch=poweroff

