# tests/auto/InvalidValue/2.31/good/Login_HoldoffTimeoutSec_infinity.target
[Unit]
Description=Foo
[Login]
HoldoffTimeoutSec=infinity

