# tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_13months.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=13months

