# tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitch_hibernate.target
[Unit]
Description=Foo
[Login]
HandleLidSwitch=hibernate

