# tests/auto/InvalidValue/2.31/good/Login_RemoveIPC_off.target
[Unit]
Description=Foo
[Login]
RemoveIPC=off

