# tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_13months.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=13months

