# tests/auto/InvalidValue/2.31/good/Login_RemoveIPC_on.target
[Unit]
Description=Foo
[Login]
RemoveIPC=on

