# tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitch_ignore.target
[Unit]
Description=Foo
[Login]
HandleLidSwitch=ignore

