# tests/auto/InvalidValue/2.31/good/Login_KillUserProcesses_1.target
[Unit]
Description=Foo
[Login]
KillUserProcesses=1

