# tests/auto/InvalidValue/2.31/good/Login_HandlePowerKey_suspend-then-hibernate.target
[Unit]
Description=Foo
[Login]
HandlePowerKey=suspend-then-hibernate

