# tests/auto/InvalidValue/2.31/good/Login_RemoveIPC_false.target
[Unit]
Description=Foo
[Login]
RemoveIPC=false

