# tests/auto/InvalidValue/2.31/good/Login_SuspendKeyIgnoreInhibited_0.target
[Unit]
Description=Foo
[Login]
SuspendKeyIgnoreInhibited=0

