# tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_14msec.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=14msec

