# tests/auto/InvalidValue/2.31/good/Login_HandleSuspendKey_kexec.target
[Unit]
Description=Foo
[Login]
HandleSuspendKey=kexec

