# tests/auto/InvalidValue/2.31/good/Login_UserStopDelaySec_14msec.target
[Unit]
Description=Foo
[Login]
UserStopDelaySec=14msec

