# tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitchExternalPower_suspend-then-hibernate.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchExternalPower=suspend-then-hibernate

