# tests/auto/InvalidValue/2.31/good/Login_IdleActionSec_22weeks.target
[Unit]
Description=Foo
[Login]
IdleActionSec=22weeks

