# tests/auto/InvalidValue/2.31/good/Login_HandleLidSwitchDocked_suspend.target
[Unit]
Description=Foo
[Login]
HandleLidSwitchDocked=suspend

