# tests/auto/InvalidValue/2.31/good/Login_InhibitDelayMaxSec_11m.target
[Unit]
Description=Foo
[Login]
InhibitDelayMaxSec=11m

