# tests/auto/InvalidValue/2.31/good/MACVTAP_Mode_vepa.target
[Unit]
Description=Foo
[MACVTAP]
Mode=vepa

