# tests/auto/InvalidValue/2.31/good/MACVTAP_Mode_private.target
[Unit]
Description=Foo
[MACVTAP]
Mode=private

