# tests/auto/InvalidValue/2.31/good/MACVLAN_Mode_bridge.target
[Unit]
Description=Foo
[MACVLAN]
Mode=bridge

