# tests/auto/InvalidValue/2.31/good/MACVLAN_Mode_passthru.target
[Unit]
Description=Foo
[MACVLAN]
Mode=passthru

