# tests/auto/InvalidValue/2.31/good/MACVTAP_Mode_passthru.target
[Unit]
Description=Foo
[MACVTAP]
Mode=passthru

