# tests/auto/InvalidValue/2.31/good/MACVLAN_Mode_vepa.target
[Unit]
Description=Foo
[MACVLAN]
Mode=vepa

