# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_14msec.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=14msec

