# tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_17second.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=17second

