# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_0days.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=0days

