# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimerAccuracySec_21us.target
[Unit]
Description=Foo
[Manager]
DefaultTimerAccuracySec=21us

