# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_9min.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=9min

