# tests/auto/InvalidValue/2.31/good/Manager_DefaultBlockIOAccounting_0.target
[Unit]
Description=Foo
[Manager]
DefaultBlockIOAccounting=0

