# tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_7minutes.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=7minutes

