# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_6h.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=6h

