# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_24w.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=24w

