# tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_13months.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=13months

