# tests/auto/InvalidValue/2.31/good/Manager_DefaultIOAccounting_1.target
[Unit]
Description=Foo
[Manager]
DefaultIOAccounting=1

