# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_13months.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=13months

