# tests/auto/InvalidValue/2.31/good/Manager_DefaultLimitSTACK_abc.target
[Unit]
Description=Foo
[Manager]
DefaultLimitSTACK=abc

