# tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_26year.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=26year

