# tests/auto/InvalidValue/2.31/good/Manager_DefaultIPAccounting_1.target
[Unit]
Description=Foo
[Manager]
DefaultIPAccounting=1

