# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimerAccuracySec_8minute.target
[Unit]
Description=Foo
[Manager]
DefaultTimerAccuracySec=8minute

