# tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_5hr.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=5hr

