# tests/auto/InvalidValue/2.31/good/Manager_DefaultBlockIOAccounting_no.target
[Unit]
Description=Foo
[Manager]
DefaultBlockIOAccounting=no

