# tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_7minutes.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=7minutes

