# tests/auto/InvalidValue/2.31/good/Manager_CrashChangeVT_no.target
[Unit]
Description=Foo
[Manager]
CrashChangeVT=no

