# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_22weeks.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=22weeks

