# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimerAccuracySec_6h.target
[Unit]
Description=Foo
[Manager]
DefaultTimerAccuracySec=6h

