# tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_14msec.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=14msec

