# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_8minute.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=8minute

