# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_26year.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=26year

