# tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_24w.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=24w

