# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStartSec_5hr.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStartSec=5hr

