# tests/auto/InvalidValue/2.31/good/Manager_DefaultLimitDATA_abc.target
[Unit]
Description=Foo
[Manager]
DefaultLimitDATA=abc

