# tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_9min.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=9min

