# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_5hr.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=5hr

