# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_3hours.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=3hours

