# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_27y.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=27y

