# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStartSec_7minutes.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStartSec=7minutes

