# tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_1day.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=1day

