# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_4hour.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=4hour

