# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_25years.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=25years

