# tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_11m.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=11m

