# tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_4hour.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=4hour

