# tests/auto/InvalidValue/2.31/good/Manager_CrashReboot_on.target
[Unit]
Description=Foo
[Manager]
CrashReboot=on

