# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimerAccuracySec_13months.target
[Unit]
Description=Foo
[Manager]
DefaultTimerAccuracySec=13months

