# tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_21us.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=21us

