# tests/auto/InvalidValue/2.31/good/Manager_CrashReboot_off.target
[Unit]
Description=Foo
[Manager]
CrashReboot=off

