# tests/auto/InvalidValue/2.31/good/Manager_DefaultMemoryAccounting_yes.target
[Unit]
Description=Foo
[Manager]
DefaultMemoryAccounting=yes

