# tests/auto/InvalidValue/2.31/good/Manager_CrashChangeVT_off.target
[Unit]
Description=Foo
[Manager]
CrashChangeVT=off

