# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_20usec.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=20usec

