# tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_5hr.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=5hr

