# tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_20usec.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=20usec

