# tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_11m.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=11m

