# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStartSec_10M.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStartSec=10M

