# tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_10M.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=10M

