# tests/auto/InvalidValue/2.31/good/Manager_DumpCore_true.target
[Unit]
Description=Foo
[Manager]
DumpCore=true

