# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_15ms.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=15ms

