# tests/auto/InvalidValue/2.31/good/Manager_DefaultMemoryAccounting_no.target
[Unit]
Description=Foo
[Manager]
DefaultMemoryAccounting=no

