# tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_6h.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=6h

