# tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_5hr.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=5hr

