# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_17second.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=17second

