# tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_27y.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=27y

