# tests/auto/InvalidValue/2.31/good/Manager_ShowStatus_on.target
[Unit]
Description=Foo
[Manager]
ShowStatus=on

