# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_8minute.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=8minute

