# tests/auto/InvalidValue/2.31/good/Manager_CrashReboot_yes.target
[Unit]
Description=Foo
[Manager]
CrashReboot=yes

