# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStartSec_4hour.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStartSec=4hour

