# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_9min.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=9min

