# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStartSec_13months.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStartSec=13months

