# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_1day.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=1day

