# tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_3hours.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=3hours

