# tests/auto/InvalidValue/2.31/good/Manager_DefaultIPAccounting_on.target
[Unit]
Description=Foo
[Manager]
DefaultIPAccounting=on

