# tests/auto/InvalidValue/2.31/good/Manager_DefaultStandardOutput_abc.target
[Unit]
Description=Foo
[Manager]
DefaultStandardOutput=abc

