# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_25years.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=25years

