# tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_20usec.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=20usec

