# tests/auto/InvalidValue/2.31/good/Manager_CrashChangeVT_0.target
[Unit]
Description=Foo
[Manager]
CrashChangeVT=0

