# tests/auto/InvalidValue/2.31/good/Manager_CrashChangeVT_1.target
[Unit]
Description=Foo
[Manager]
CrashChangeVT=1

