# tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_4hour.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=4hour

