# tests/auto/InvalidValue/2.31/good/Manager_CrashReboot_false.target
[Unit]
Description=Foo
[Manager]
CrashReboot=false

