# tests/auto/InvalidValue/2.31/good/Manager_DefaultBlockIOAccounting_off.target
[Unit]
Description=Foo
[Manager]
DefaultBlockIOAccounting=off

