# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_24w.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=24w

