# tests/auto/InvalidValue/2.31/good/Manager_ShowStatus_yes.target
[Unit]
Description=Foo
[Manager]
ShowStatus=yes

