# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStartSec_11m.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStartSec=11m

