# tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_26year.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=26year

