# tests/auto/InvalidValue/2.31/good/Manager_AmbientCapabilities_abc.target
[Unit]
Description=Foo
[Manager]
AmbientCapabilities=abc

