# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStartSec_27y.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStartSec=27y

