# tests/auto/InvalidValue/2.31/good/Manager_DefaultCPUAccounting_on.target
[Unit]
Description=Foo
[Manager]
DefaultCPUAccounting=on

