# tests/auto/InvalidValue/2.31/good/Manager_DefaultTasksMax_9999999999999%25.target
[Unit]
Description=Foo
[Manager]
DefaultTasksMax=9999999999999%

