# tests/auto/InvalidValue/2.31/good/Manager_DefaultTasksAccounting_no.target
[Unit]
Description=Foo
[Manager]
DefaultTasksAccounting=no

