# tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_5hr.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=5hr

