# tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_14msec.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=14msec

