# tests/auto/InvalidValue/2.31/good/Manager_LogColor_0.target
[Unit]
Description=Foo
[Manager]
LogColor=0

