# tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_25years.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=25years

