# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimerAccuracySec_15ms.target
[Unit]
Description=Foo
[Manager]
DefaultTimerAccuracySec=15ms

