# tests/auto/InvalidValue/2.31/good/Manager_LogColor_yes.target
[Unit]
Description=Foo
[Manager]
LogColor=yes

