# tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_infinity.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=infinity

