# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_9min.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=9min

