# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_26year.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=26year

