# tests/auto/InvalidValue/2.31/good/Manager_CrashReboot_true.target
[Unit]
Description=Foo
[Manager]
CrashReboot=true

