# tests/auto/InvalidValue/2.31/good/Manager_DefaultStandardError_abc.target
[Unit]
Description=Foo
[Manager]
DefaultStandardError=abc

