# tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_11m.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=11m

