# tests/auto/InvalidValue/2.31/good/Manager_CrashShell_0.target
[Unit]
Description=Foo
[Manager]
CrashShell=0

