# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStartSec_3hours.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStartSec=3hours

