# tests/auto/InvalidValue/2.31/good/Manager_DefaultLimitNPROC_abc.target
[Unit]
Description=Foo
[Manager]
DefaultLimitNPROC=abc

