# tests/auto/InvalidValue/2.31/good/Manager_LogColor_off.target
[Unit]
Description=Foo
[Manager]
LogColor=off

