# tests/auto/InvalidValue/2.31/good/Manager_DefaultCPUAccounting_no.target
[Unit]
Description=Foo
[Manager]
DefaultCPUAccounting=no

