# tests/auto/InvalidValue/2.31/good/Manager_DumpCore_false.target
[Unit]
Description=Foo
[Manager]
DumpCore=false

