# tests/auto/InvalidValue/2.31/good/Manager_LogLocation_abc.target
[Unit]
Description=Foo
[Manager]
LogLocation=abc

