# tests/auto/InvalidValue/2.31/good/Manager_DefaultMemoryAccounting_false.target
[Unit]
Description=Foo
[Manager]
DefaultMemoryAccounting=false

