# tests/auto/InvalidValue/2.31/good/Manager_LogColor_no.target
[Unit]
Description=Foo
[Manager]
LogColor=no

