# tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_6h.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=6h

