# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_12month.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=12month

