# tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_7minutes.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=7minutes

