# tests/auto/InvalidValue/2.31/good/Manager_DefaultIPAccounting_false.target
[Unit]
Description=Foo
[Manager]
DefaultIPAccounting=false

