# tests/auto/InvalidValue/2.31/good/Manager_DumpCore_yes.target
[Unit]
Description=Foo
[Manager]
DumpCore=yes

