# tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_23week.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=23week

