# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStartSec_28%C2%B5s.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStartSec=28µs

