# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimerAccuracySec_2d.target
[Unit]
Description=Foo
[Manager]
DefaultTimerAccuracySec=2d

