# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_13months.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=13months

