# tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_23week.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=23week

