# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_15ms.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=15ms

