# tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_2d.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=2d

