# tests/auto/InvalidValue/2.31/good/Manager_ShutdownWatchdogSec_13months.target
[Unit]
Description=Foo
[Manager]
ShutdownWatchdogSec=13months

