# tests/auto/InvalidValue/2.31/good/Manager_LogTarget_abc.target
[Unit]
Description=Foo
[Manager]
LogTarget=abc

