# tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_19s.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=19s

