# tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_13months.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=13months

