# tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_3hours.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=3hours

