# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_13months.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=13months

