# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitIntervalSec_14msec.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitIntervalSec=14msec

