# tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_11m.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=11m

