# tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_17second.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=17second

