# tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_16seconds.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=16seconds

