# tests/auto/InvalidValue/2.31/good/Manager_DefaultRestartSec_22weeks.target
[Unit]
Description=Foo
[Manager]
DefaultRestartSec=22weeks

