# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_27y.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=27y

