# tests/auto/InvalidValue/2.31/good/Manager_CrashShell_no.target
[Unit]
Description=Foo
[Manager]
CrashShell=no

