# tests/auto/InvalidValue/2.31/good/Manager_DefaultBlockIOAccounting_yes.target
[Unit]
Description=Foo
[Manager]
DefaultBlockIOAccounting=yes

