# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_20usec.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=20usec

