# tests/auto/InvalidValue/2.31/good/Manager_SystemCallArchitectures_27y.target
[Unit]
Description=Foo
[Manager]
SystemCallArchitectures=27y

