# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStopSec_16seconds.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStopSec=16seconds

