# tests/auto/InvalidValue/2.31/good/Manager_CapabilityBoundingSet_abc.target
[Unit]
Description=Foo
[Manager]
CapabilityBoundingSet=abc

