# tests/auto/InvalidValue/2.31/good/Manager_DefaultLimitAS_abc.target
[Unit]
Description=Foo
[Manager]
DefaultLimitAS=abc

