# tests/auto/InvalidValue/2.31/good/Manager_DefaultEnvironment_a%3Db.target
[Unit]
Description=Foo
[Manager]
DefaultEnvironment=a=b

