# tests/auto/InvalidValue/2.31/good/Manager_DefaultMemoryAccounting_true.target
[Unit]
Description=Foo
[Manager]
DefaultMemoryAccounting=true

