# tests/auto/InvalidValue/2.31/good/Manager_DumpCore_1.target
[Unit]
Description=Foo
[Manager]
DumpCore=1

