# tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_9min.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=9min

