# tests/auto/InvalidValue/2.31/good/Manager_RuntimeWatchdogSec_13months.target
[Unit]
Description=Foo
[Manager]
RuntimeWatchdogSec=13months

