# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_23week.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=23week

