# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimerAccuracySec_infinity.target
[Unit]
Description=Foo
[Manager]
DefaultTimerAccuracySec=infinity

