# tests/auto/InvalidValue/2.31/good/Manager_DefaultMemoryAccounting_0.target
[Unit]
Description=Foo
[Manager]
DefaultMemoryAccounting=0

