# tests/auto/InvalidValue/2.31/good/Manager_DefaultTimeoutStartSec_26year.target
[Unit]
Description=Foo
[Manager]
DefaultTimeoutStartSec=26year

