# tests/auto/InvalidValue/2.31/good/Manager_ShowStatus_off.target
[Unit]
Description=Foo
[Manager]
ShowStatus=off

