# tests/auto/InvalidValue/2.31/good/Manager_DefaultStartLimitBurst_17second.target
[Unit]
Description=Foo
[Manager]
DefaultStartLimitBurst=17second

