# tests/auto/InvalidValue/2.31/good/Manager_DefaultMemoryAccounting_1.target
[Unit]
Description=Foo
[Manager]
DefaultMemoryAccounting=1

