# tests/auto/InvalidValue/2.31/good/Manager_DefaultLimitMSGQUEUE_abc.target
[Unit]
Description=Foo
[Manager]
DefaultLimitMSGQUEUE=abc

