# tests/auto/InvalidValue/2.31/good/Manager_CrashReboot_1.target
[Unit]
Description=Foo
[Manager]
CrashReboot=1

