# tests/auto/InvalidValue/2.31/good/Manager_DefaultLimitFSIZE_abc.target
[Unit]
Description=Foo
[Manager]
DefaultLimitFSIZE=abc

