# tests/auto/InvalidValue/2.31/good/Manager_TimerSlackNSec_15ms.target
[Unit]
Description=Foo
[Manager]
TimerSlackNSec=15ms

