# tests/auto/InvalidValue/2.31/good/Match_Name_abc.target
[Unit]
Description=Foo
[Match]
Name=abc

