# tests/auto/InvalidValue/2.31/good/Match_Driver_abc.target
[Unit]
Description=Foo
[Match]
Driver=abc

