# tests/auto/InvalidValue/2.31/good/Match_Path_abc.target
[Unit]
Description=Foo
[Match]
Path=abc

