# tests/auto/InvalidValue/2.31/good/Match_Host_abc.target
[Unit]
Description=Foo
[Match]
Host=abc

