# tests/auto/InvalidValue/2.31/good/Mount_RootImage_--usr--sbin.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
RootImage=/usr/sbin

