# tests/auto/InvalidValue/2.31/good/Mount_TimerSlackNSec_14msec.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
TimerSlackNSec=14msec

