# tests/auto/InvalidValue/2.31/good/Mount_CPUQuota_0%25.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
CPUQuota=0%

