# tests/auto/InvalidValue/2.31/good/Mount_FinalKillSignal_SIGTSTP.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
FinalKillSignal=SIGTSTP

