# tests/auto/InvalidValue/2.31/good/Mount_SloppyOptions_1.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
SloppyOptions=1

