# tests/auto/InvalidValue/2.31/good/Mount_BlockIOReadBandwidthMax_--dev--bar%200B.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
BlockIOReadBandwidthMax=/dev/bar 0B

