# tests/auto/InvalidValue/2.31/good/Mount_IOWriteBandwidthMax_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
IOWriteBandwidthMax=abc

