# tests/auto/InvalidValue/2.31/good/Mount_StartupCPUWeight_1.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
StartupCPUWeight=1

