# tests/auto/InvalidValue/2.31/good/Mount_SecureBits_no-setuid-fixup-locked.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
SecureBits=no-setuid-fixup-locked

