# tests/auto/InvalidValue/2.31/good/Mount_Personality_ppc64-le.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
Personality=ppc64-le

