# tests/auto/InvalidValue/2.31/good/Mount_ProtectHome_read-only.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
ProtectHome=read-only

