# tests/auto/InvalidValue/2.31/good/Mount_TimeoutSec_infinity.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
TimeoutSec=infinity

