# tests/auto/InvalidValue/2.31/good/Mount_User_metoo.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
User=metoo

