# tests/auto/InvalidValue/2.31/good/Mount_CPUSchedulingPriority_99.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
CPUSchedulingPriority=99

