# tests/auto/InvalidValue/2.31/good/Mount_LimitRTPRIO_0.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
LimitRTPRIO=0

