# tests/auto/InvalidValue/2.31/good/Mount_LimitLOCKS_infinity.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
LimitLOCKS=infinity

