# tests/auto/InvalidValue/2.31/good/Mount_BlockIOWriteBandwidthMax_--dev--bar%209999999999999B.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
BlockIOWriteBandwidthMax=/dev/bar 9999999999999B

