# tests/auto/InvalidValue/2.31/good/Mount_StartupCPUShares_2.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
StartupCPUShares=2

