# tests/auto/InvalidValue/2.31/good/Mount_KillSignal_SIGSTOP.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
KillSignal=SIGSTOP

