# tests/auto/InvalidValue/2.31/good/Mount_BlockIODeviceWeight_--dev--bar%201.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
BlockIODeviceWeight=/dev/bar 1

