# tests/auto/InvalidValue/2.31/good/NetDev_Name_abc.target
[Unit]
Description=Foo
[NetDev]
Name=abc
Name=abc

