# tests/auto/InvalidValue/2.31/good/NetDev_MTUBytes_9999999999999K.target
[Unit]
Description=Foo
[NetDev]
Name=abc
MTUBytes=9999999999999K

