# tests/auto/InvalidValue/2.31/good/NetDev_Description_abc.target
[Unit]
Description=Foo
[NetDev]
Name=abc
Description=abc

