# tests/auto/InvalidValue/2.31/good/NetDev_MACAddress_01%3A23%3A45%3A67%3A89%3Aab.target
[Unit]
Description=Foo
[NetDev]
Name=abc
MACAddress=01:23:45:67:89:ab

