# tests/auto/InvalidValue/2.31/good/NetDev_Kind_geneve.target
[Unit]
Description=Foo
[NetDev]
Name=abc
Kind=geneve

