# tests/auto/InvalidValue/2.31/good/Network_PrimarySlave_off.target
[Unit]
Description=Foo
[Network]
PrimarySlave=off

