# tests/auto/InvalidValue/2.31/good/Network_LinkLocalAddressing_ipv6.target
[Unit]
Description=Foo
[Network]
LinkLocalAddressing=ipv6

