# tests/auto/InvalidValue/2.31/good/Network_DNSSEC_no.target
[Unit]
Description=Foo
[Network]
DNSSEC=no

