# tests/auto/InvalidValue/2.31/good/Network_VLAN_abc.target
[Unit]
Description=Foo
[Network]
VLAN=abc

