# tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_9min.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=9min

