# tests/auto/InvalidValue/2.31/good/Network_DNSDefaultRoute_0.target
[Unit]
Description=Foo
[Network]
DNSDefaultRoute=0

