# tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_13months.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=13months

