# tests/auto/InvalidValue/2.31/good/Network_ActiveSlave_1.target
[Unit]
Description=Foo
[Network]
ActiveSlave=1

