# tests/auto/InvalidValue/2.31/good/Network_IPv6ProxyNDPAddress_1.1.1.1.target
[Unit]
Description=Foo
[Network]
IPv6ProxyNDPAddress=1.1.1.1

