# tests/auto/InvalidValue/2.31/good/Network_DHCPServer_true.target
[Unit]
Description=Foo
[Network]
DHCPServer=true

