# tests/auto/InvalidValue/2.31/good/Network_PrimarySlave_0.target
[Unit]
Description=Foo
[Network]
PrimarySlave=0

