# tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_10M.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=10M

