# tests/auto/InvalidValue/2.31/good/Network_LLDP_yes.target
[Unit]
Description=Foo
[Network]
LLDP=yes

