# tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_5hr.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=5hr

