# tests/auto/InvalidValue/2.31/good/Network_DNSDefaultRoute_on.target
[Unit]
Description=Foo
[Network]
DNSDefaultRoute=on

