# tests/auto/InvalidValue/2.31/good/Network_LLMNR_no.target
[Unit]
Description=Foo
[Network]
LLMNR=no

