# tests/auto/InvalidValue/2.31/good/Network_EmitLLDP_yes.target
[Unit]
Description=Foo
[Network]
EmitLLDP=yes

