# tests/auto/InvalidValue/2.31/good/Network_MulticastDNS_1.target
[Unit]
Description=Foo
[Network]
MulticastDNS=1

