# tests/auto/InvalidValue/2.31/good/Network_EmitLLDP_false.target
[Unit]
Description=Foo
[Network]
EmitLLDP=false

