# tests/auto/InvalidValue/2.31/good/Network_IPv6MTUBytes_9999999999999.target
[Unit]
Description=Foo
[Network]
IPv6MTUBytes=9999999999999

