# tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_2d.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=2d

