# tests/auto/InvalidValue/2.31/good/Network_DHCPServer_on.target
[Unit]
Description=Foo
[Network]
DHCPServer=on

