# tests/auto/InvalidValue/2.31/good/Network_DNSDefaultRoute_true.target
[Unit]
Description=Foo
[Network]
DNSDefaultRoute=true

