# tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_27y.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=27y

