# tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_22weeks.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=22weeks

