# tests/auto/InvalidValue/2.31/good/Network_LinkLocalAddressing_no.target
[Unit]
Description=Foo
[Network]
LinkLocalAddressing=no

