# tests/auto/InvalidValue/2.31/good/Network_PrimarySlave_no.target
[Unit]
Description=Foo
[Network]
PrimarySlave=no

