# tests/auto/InvalidValue/2.31/good/Network_ActiveSlave_off.target
[Unit]
Description=Foo
[Network]
ActiveSlave=off

