# tests/auto/InvalidValue/2.31/good/Network_EmitLLDP_non-tpmr-bridge.target
[Unit]
Description=Foo
[Network]
EmitLLDP=non-tpmr-bridge

