# tests/auto/InvalidValue/2.31/good/Network_LLDP_true.target
[Unit]
Description=Foo
[Network]
LLDP=true

