# tests/auto/InvalidValue/2.31/good/Network_LinkLocalAddressing_ipv4.target
[Unit]
Description=Foo
[Network]
LinkLocalAddressing=ipv4

