# tests/auto/InvalidValue/2.31/good/Network_DNSSECNegativeTrustAnchors_abc.target
[Unit]
Description=Foo
[Network]
DNSSECNegativeTrustAnchors=abc

