# tests/auto/InvalidValue/2.31/good/Network_EmitLLDP_no.target
[Unit]
Description=Foo
[Network]
EmitLLDP=no

