# tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_20usec.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=20usec

