# tests/auto/InvalidValue/2.31/good/Network_DHCPServer_0.target
[Unit]
Description=Foo
[Network]
DHCPServer=0

