# tests/auto/InvalidValue/2.31/good/Network_ConfigureWithoutCarrier_false.target
[Unit]
Description=Foo
[Network]
ConfigureWithoutCarrier=false

