# tests/auto/InvalidValue/2.31/good/Network_DNSDefaultRoute_no.target
[Unit]
Description=Foo
[Network]
DNSDefaultRoute=no

