# tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_false.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=false

