# tests/auto/InvalidValue/2.31/good/Network_DNSSEC_yes.target
[Unit]
Description=Foo
[Network]
DNSSEC=yes

