# tests/auto/InvalidValue/2.31/good/Network_NTP_255.255.255.255--32.target
[Unit]
Description=Foo
[Network]
NTP=255.255.255.255/32

