# tests/auto/InvalidValue/2.31/good/Network_Tunnel_abc.target
[Unit]
Description=Foo
[Network]
Tunnel=abc

