# tests/auto/InvalidValue/2.31/good/Network_DHCPServer_false.target
[Unit]
Description=Foo
[Network]
DHCPServer=false

