# tests/auto/InvalidValue/2.31/good/Network_Gateway__ipv6ra.target
[Unit]
Description=Foo
[Network]
Gateway=_ipv6ra

