# tests/auto/InvalidValue/2.31/good/Network_EmitLLDP_0.target
[Unit]
Description=Foo
[Network]
EmitLLDP=0

