# tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_1day.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=1day

