# tests/auto/InvalidValue/2.31/good/Network_IPv6HopLimit_9999999999999.target
[Unit]
Description=Foo
[Network]
IPv6HopLimit=9999999999999

