# tests/auto/InvalidValue/2.31/good/Network_LLMNR_0.target
[Unit]
Description=Foo
[Network]
LLMNR=0

