# tests/auto/InvalidValue/2.31/good/Network_Gateway__dhcp4.target
[Unit]
Description=Foo
[Network]
Gateway=_dhcp4

