# tests/auto/InvalidValue/2.31/good/Network_BindCarrier_abc.target
[Unit]
Description=Foo
[Network]
BindCarrier=abc

