# tests/auto/InvalidValue/2.31/good/Network_VRF_abc.target
[Unit]
Description=Foo
[Network]
VRF=abc

