# tests/auto/InvalidValue/2.31/good/Network_LLMNR_1.target
[Unit]
Description=Foo
[Network]
LLMNR=1

