# tests/auto/InvalidValue/2.31/good/Network_DHCPServer_yes.target
[Unit]
Description=Foo
[Network]
DHCPServer=yes

