# tests/auto/InvalidValue/2.31/good/Network_DHCPServer_off.target
[Unit]
Description=Foo
[Network]
DHCPServer=off

