# tests/auto/InvalidValue/2.31/good/Network_DNSOverTLS_false.target
[Unit]
Description=Foo
[Network]
DNSOverTLS=false

