# tests/auto/InvalidValue/2.31/good/Network_LLMNR_true.target
[Unit]
Description=Foo
[Network]
LLMNR=true

