# tests/auto/InvalidValue/2.31/good/Network_MACVLAN_abc.target
[Unit]
Description=Foo
[Network]
MACVLAN=abc

