# tests/auto/InvalidValue/2.31/good/Network_Bridge_abc.target
[Unit]
Description=Foo
[Network]
Bridge=abc

