# tests/auto/InvalidValue/2.31/good/Network_ActiveSlave_no.target
[Unit]
Description=Foo
[Network]
ActiveSlave=no

