# tests/auto/InvalidValue/2.31/good/Network_VXLAN_abc.target
[Unit]
Description=Foo
[Network]
VXLAN=abc

