# tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_7minutes.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=7minutes

