# tests/auto/InvalidValue/2.31/good/Network_LLMNR_yes.target
[Unit]
Description=Foo
[Network]
LLMNR=yes

