# tests/auto/InvalidValue/2.31/good/Network_IPv6PrivacyExtensions_kernel.target
[Unit]
Description=Foo
[Network]
IPv6PrivacyExtensions=kernel

