# tests/auto/InvalidValue/2.31/good/Network_Domains_abc.target
[Unit]
Description=Foo
[Network]
Domains=abc

