# tests/auto/InvalidValue/2.31/good/Network_MulticastDNS_true.target
[Unit]
Description=Foo
[Network]
MulticastDNS=true

