# tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_28%C2%B5s.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=28µs

