# tests/auto/InvalidValue/2.31/good/Network_LLDP_no.target
[Unit]
Description=Foo
[Network]
LLDP=no

