# tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_3hours.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=3hours

