# tests/auto/InvalidValue/2.31/good/Network_IPv4LLRoute_yes.target
[Unit]
Description=Foo
[Network]
IPv4LLRoute=yes

