# tests/auto/InvalidValue/2.31/good/Network_PrimarySlave_false.target
[Unit]
Description=Foo
[Network]
PrimarySlave=false

