# tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_15ms.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=15ms

