# tests/auto/InvalidValue/2.31/good/Network_EmitLLDP_true.target
[Unit]
Description=Foo
[Network]
EmitLLDP=true

