# tests/auto/InvalidValue/2.31/good/Network_IgnoreCarrierLoss_16seconds.target
[Unit]
Description=Foo
[Network]
IgnoreCarrierLoss=16seconds

