# tests/auto/InvalidValue/2.31/good/Network_DHCPServer_1.target
[Unit]
Description=Foo
[Network]
DHCPServer=1

