# tests/auto/InvalidValue/2.31/good/Network_PrimarySlave_true.target
[Unit]
Description=Foo
[Network]
PrimarySlave=true

