# tests/auto/InvalidValue/2.31/good/Network_DNS_255.255.255.255--32.target
[Unit]
Description=Foo
[Network]
DNS=255.255.255.255/32

