# tests/auto/InvalidValue/2.31/good/Path_DirectoryNotEmpty_--dev--.target
[Unit]
Description=Foo
[Path]
DirectoryNotEmpty=/dev/

