# tests/auto/InvalidValue/2.31/good/Path_PathChanged_--dev--.target
[Unit]
Description=Foo
[Path]
PathChanged=/dev/

