# tests/auto/InvalidValue/2.31/good/Path_PathExists_--tmp.target
[Unit]
Description=Foo
[Path]
PathExists=/tmp

