# tests/auto/InvalidValue/2.31/good/Path_PathModified_--tmp.target
[Unit]
Description=Foo
[Path]
PathModified=/tmp

