# tests/auto/InvalidValue/2.31/good/Path_PathExistsGlob_abc.target
[Unit]
Description=Foo
[Path]
PathExistsGlob=abc

