# tests/auto/InvalidValue/2.31/good/Path_MakeDirectory_false.target
[Unit]
Description=Foo
[Path]
MakeDirectory=false

