# tests/auto/InvalidValue/2.31/good/Path_MakeDirectory_0.target
[Unit]
Description=Foo
[Path]
MakeDirectory=0

