# tests/auto/InvalidValue/2.31/good/Path_DirectoryMode_004.target
[Unit]
Description=Foo
[Path]
DirectoryMode=004

