# tests/auto/InvalidValue/2.31/good/Peer_MACAddress_00-11-22-33-44-55.target
[Unit]
Description=Foo
[Peer]
Name=abc
MACAddress=00-11-22-33-44-55

