# tests/auto/InvalidValue/2.31/good/Route_Protocol_255.target
[Unit]
Description=Foo
[Route]
Protocol=255

