# tests/auto/InvalidValue/2.31/good/Route_Metric_9999999999999.target
[Unit]
Description=Foo
[Route]
Metric=9999999999999

