# tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_To_1.1.1.1--0.target
[Unit]
Description=Foo
[RoutingPolicyRule]
To=1.1.1.1/0

