# tests/auto/InvalidValue/2.31/good/Route_QuickAck_false.target
[Unit]
Description=Foo
[Route]
QuickAck=false

