# tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_From_255.255.255.255--32.target
[Unit]
Description=Foo
[RoutingPolicyRule]
From=255.255.255.255/32

