# tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_Priority_0.target
[Unit]
Description=Foo
[RoutingPolicyRule]
Priority=0

