# tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_InvertRule_off.target
[Unit]
Description=Foo
[RoutingPolicyRule]
InvertRule=off

