# tests/auto/InvalidValue/2.31/good/Route_Destination_abc.target
[Unit]
Description=Foo
[Route]
Destination=abc

