# tests/auto/InvalidValue/2.31/good/Route_GatewayOnLink_true.target
[Unit]
Description=Foo
[Route]
GatewayOnLink=true

