# tests/auto/InvalidValue/2.31/good/Route_GatewayOnLink_yes.target
[Unit]
Description=Foo
[Route]
GatewayOnLink=yes

