# tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_IPProtocol_tcp.target
[Unit]
Description=Foo
[RoutingPolicyRule]
IPProtocol=tcp

