# tests/auto/InvalidValue/2.31/good/Route_Protocol_boot.target
[Unit]
Description=Foo
[Route]
Protocol=boot

