# tests/auto/InvalidValue/2.31/good/Route_Source_abc.target
[Unit]
Description=Foo
[Route]
Source=abc

