# tests/auto/InvalidValue/2.31/good/Route_Protocol_kernel.target
[Unit]
Description=Foo
[Route]
Protocol=kernel

