# tests/auto/InvalidValue/2.31/good/Route_Type_throw.target
[Unit]
Description=Foo
[Route]
Type=throw

