# tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_FirewallMark_1.target
[Unit]
Description=Foo
[RoutingPolicyRule]
FirewallMark=1

