# tests/auto/InvalidValue/2.31/good/Route_Scope_global.target
[Unit]
Description=Foo
[Route]
Scope=global

