# tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_InvertRule_false.target
[Unit]
Description=Foo
[RoutingPolicyRule]
InvertRule=false

