# tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_InvertRule_0.target
[Unit]
Description=Foo
[RoutingPolicyRule]
InvertRule=0

