# tests/auto/InvalidValue/2.31/good/Route_IPv6Preference_high.target
[Unit]
Description=Foo
[Route]
IPv6Preference=high

