# tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_InvertRule_no.target
[Unit]
Description=Foo
[RoutingPolicyRule]
InvertRule=no

