# tests/auto/InvalidValue/2.31/good/Route_InitialAdvertisedReceiveWindow_4294967295.target
[Unit]
Description=Foo
[Route]
InitialAdvertisedReceiveWindow=4294967295

