# tests/auto/InvalidValue/2.31/good/Route_PreferredSource_abc.target
[Unit]
Description=Foo
[Route]
PreferredSource=abc

