# tests/auto/InvalidValue/2.31/good/Route_Type_xresolve.target
[Unit]
Description=Foo
[Route]
Type=xresolve

