# tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_TypeOfService_0.target
[Unit]
Description=Foo
[RoutingPolicyRule]
TypeOfService=0

