# tests/auto/InvalidValue/2.31/good/Route_Metric_0.target
[Unit]
Description=Foo
[Route]
Metric=0

