# tests/auto/InvalidValue/2.31/good/RoutingPolicyRule_DestinationPort_0%200.target
[Unit]
Description=Foo
[RoutingPolicyRule]
DestinationPort=0 0

