# tests/auto/InvalidValue/2.31/good/Route_GatewayOnLink_false.target
[Unit]
Description=Foo
[Route]
GatewayOnLink=false

