# tests/auto/InvalidValue/2.31/good/Route_QuickAck_1.target
[Unit]
Description=Foo
[Route]
QuickAck=1

