# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_14.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=14

