# tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_21us.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=21us

