# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGCONT.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGCONT

