# tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_11m.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=11m

