# tests/auto/InvalidValue/2.31/good/Service_FileDescriptorStoreMax_9999999999999.target
[Unit]
Description=Foo
[Service]
FileDescriptorStoreMax=9999999999999

