# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ESHUTDOWN.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ESHUTDOWN

