# tests/auto/InvalidValue/2.31/good/Service_RestartSec_17second.target
[Unit]
Description=Foo
[Service]
RestartSec=17second

