# tests/auto/InvalidValue/2.31/good/Service_ProtectSystem_yes.target
[Unit]
Description=Foo
[Service]
ProtectSystem=yes

