# tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGUSR1.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGUSR1

