# tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGSYS.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGSYS

