# tests/auto/InvalidValue/2.31/good/Service_TTYReset_false.target
[Unit]
Description=Foo
[Service]
TTYReset=false

