# tests/auto/InvalidValue/2.31/good/Service_LimitMEMLOCK_infinity.target
[Unit]
Description=Foo
[Service]
LimitMEMLOCK=infinity

