# tests/auto/InvalidValue/2.31/good/Service_PrivateNetwork_false.target
[Unit]
Description=Foo
[Service]
PrivateNetwork=false

