# tests/auto/InvalidValue/2.31/good/Service_LimitSTACK_9999999999999.target
[Unit]
Description=Foo
[Service]
LimitSTACK=9999999999999

