# tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGSTKFLT.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGSTKFLT

