# tests/auto/InvalidValue/2.31/good/Service_TTYVHangup_on.target
[Unit]
Description=Foo
[Service]
TTYVHangup=on

