# tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_20usec.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=20usec

