# tests/auto/InvalidValue/2.31/good/Service_Type_idle.target
[Unit]
Description=Foo
[Service]
Type=idle

