# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGTTOU.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGTTOU

