# tests/auto/InvalidValue/2.31/good/Service_Delegate_cpuacct.target
[Unit]
Description=Foo
[Service]
Delegate=cpuacct

