# tests/auto/InvalidValue/2.31/good/Service_ExecReload_--bin--sh.target
[Unit]
Description=Foo
[Service]
ExecReload=/bin/sh

