# tests/auto/InvalidValue/2.31/good/Service_RemoveIPC_0.target
[Unit]
Description=Foo
[Service]
RemoveIPC=0

