# tests/auto/InvalidValue/2.31/good/Service_ExecStartPre_--bin--sh.target
[Unit]
Description=Foo
[Service]
ExecStartPre=/bin/sh

