# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGURG.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGURG

