# tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_12month.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=12month

