# tests/auto/InvalidValue/2.31/good/Service_KillMode_mixed.target
[Unit]
Description=Foo
[Service]
KillMode=mixed

