# tests/auto/InvalidValue/2.31/good/Service_RootDirectoryStartOnly_yes.target
[Unit]
Description=Foo
[Service]
RootDirectoryStartOnly=yes

