# tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_4hour.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=4hour

