# tests/auto/InvalidValue/2.31/good/Service_InaccessiblePaths_abc.target
[Unit]
Description=Foo
[Service]
InaccessiblePaths=abc

