# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGINT.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGINT

