# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGILL.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGILL

