# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_3hours.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=3hours

