# tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_SIGPOLL.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=SIGPOLL

