# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGTRAP.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGTRAP

