# tests/auto/InvalidValue/2.31/good/Service_Delegate_yes.target
[Unit]
Description=Foo
[Service]
Delegate=yes

