# tests/auto/InvalidValue/2.31/good/Service_NonBlocking_false.target
[Unit]
Description=Foo
[Service]
NonBlocking=false

