# tests/auto/InvalidValue/2.31/good/Service_TTYPath_--dev--foo.target
[Unit]
Description=Foo
[Service]
TTYPath=/dev/foo

