# tests/auto/InvalidValue/2.31/good/Service_LimitSTACK_infinity.target
[Unit]
Description=Foo
[Service]
LimitSTACK=infinity

