# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGCONT.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGCONT

