# tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_infinity.target
[Unit]
Description=Foo
[Service]
TimeoutSec=infinity

