# tests/auto/InvalidValue/2.31/good/Service_SendSIGKILL_yes.target
[Unit]
Description=Foo
[Service]
SendSIGKILL=yes

