# tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_15ms.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=15ms

