# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGIOT.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGIOT

