# tests/auto/InvalidValue/2.31/good/Service_TTYReset_no.target
[Unit]
Description=Foo
[Service]
TTYReset=no

