# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_52.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=52

