# tests/auto/InvalidValue/2.31/good/Service_DisableControllers_io.target
[Unit]
Description=Foo
[Service]
DisableControllers=io

