# tests/auto/InvalidValue/2.31/good/Service_LimitSIGPENDING_infinity.target
[Unit]
Description=Foo
[Service]
LimitSIGPENDING=infinity

