# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ECHRNG.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ECHRNG

