# tests/auto/InvalidValue/2.31/good/Service_SyslogFacility_local0.target
[Unit]
Description=Foo
[Service]
SyslogFacility=local0

