# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGPOLL.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGPOLL

