# tests/auto/InvalidValue/2.31/good/Service_RestrictRealtime_off.target
[Unit]
Description=Foo
[Service]
RestrictRealtime=off

