# tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_17second.target
[Unit]
Description=Foo
[Service]
WatchdogSec=17second

