# tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_infinity.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=infinity

