# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGTRAP.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGTRAP

