# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EILSEQ.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EILSEQ

