# tests/auto/InvalidValue/2.31/good/Service_RemainAfterExit_true.target
[Unit]
Description=Foo
[Service]
RemainAfterExit=true

