# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_23.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=23

