# tests/auto/InvalidValue/2.31/good/Service_PrivateDevices_true.target
[Unit]
Description=Foo
[Service]
PrivateDevices=true

