# tests/auto/InvalidValue/2.31/good/Service_SupplementaryGroups_iamagroup.target
[Unit]
Description=Foo
[Service]
SupplementaryGroups=iamagroup

