# tests/auto/InvalidValue/2.31/good/Service_RestartSec_14msec.target
[Unit]
Description=Foo
[Service]
RestartSec=14msec

