# tests/auto/InvalidValue/2.31/good/Service_IOSchedulingClass_realtime.target
[Unit]
Description=Foo
[Service]
IOSchedulingClass=realtime

