# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ECONNABORTED.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ECONNABORTED

