# tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_51.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=51

