# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGTSTP.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGTSTP

