# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGVTALRM.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGVTALRM

