# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGSTKFLT.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGSTKFLT

