# tests/auto/InvalidValue/2.31/good/Service_NotifyAccess_exec.target
[Unit]
Description=Foo
[Service]
NotifyAccess=exec

