# tests/auto/InvalidValue/2.31/good/Service_SendSIGHUP_yes.target
[Unit]
Description=Foo
[Service]
SendSIGHUP=yes

