# tests/auto/InvalidValue/2.31/good/Service_StartupIOWeight_1.target
[Unit]
Description=Foo
[Service]
StartupIOWeight=1

