# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EMSGSIZE.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EMSGSIZE

