# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_22weeks.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=22weeks

