# tests/auto/InvalidValue/2.31/good/Service_StandardError_journal%2Bconsole.target
[Unit]
Description=Foo
[Service]
StandardError=journal+console

