# tests/auto/InvalidValue/2.31/good/Service_PrivateNetwork_yes.target
[Unit]
Description=Foo
[Service]
PrivateNetwork=yes

