# tests/auto/InvalidValue/2.31/good/Service_IOWeight_1.target
[Unit]
Description=Foo
[Service]
IOWeight=1

