# tests/auto/InvalidValue/2.31/good/Service_StartLimitIntervalSec_28%C2%B5s.target
[Unit]
Description=Foo
[Service]
StartLimitIntervalSec=28µs

