# tests/auto/InvalidValue/2.31/good/Service_RemainAfterExit_no.target
[Unit]
Description=Foo
[Service]
RemainAfterExit=no

