# tests/auto/InvalidValue/2.31/good/Service_SendSIGHUP_on.target
[Unit]
Description=Foo
[Service]
SendSIGHUP=on

