# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_21us.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=21us

