# tests/auto/InvalidValue/2.31/good/Service_IgnoreSIGPIPE_1.target
[Unit]
Description=Foo
[Service]
IgnoreSIGPIPE=1

