# tests/auto/InvalidValue/2.31/good/Service_ReadWriteDirectories_abc.target
[Unit]
Description=Foo
[Service]
ReadWriteDirectories=abc

