# tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_14msec.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=14msec

