# tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_15ms.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=15ms

