# tests/auto/InvalidValue/2.31/good/Service_IOReadBandwidthMax_abc.target
[Unit]
Description=Foo
[Service]
IOReadBandwidthMax=abc

