# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGABRT.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGABRT

