# tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGCLD.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGCLD

