# tests/auto/InvalidValue/2.31/good/Service_KillSignal_38.target
[Unit]
Description=Foo
[Service]
KillSignal=38

