# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_55.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=55

