# tests/auto/InvalidValue/2.31/good/Service_LimitSIGPENDING_9999999999999.target
[Unit]
Description=Foo
[Service]
LimitSIGPENDING=9999999999999

