# tests/auto/InvalidValue/2.31/good/Service_SendSIGKILL_on.target
[Unit]
Description=Foo
[Service]
SendSIGKILL=on

