# tests/auto/InvalidValue/2.31/good/Service_BlockIOReadBandwidthMax_--dev--bar%200K.target
[Unit]
Description=Foo
[Service]
BlockIOReadBandwidthMax=/dev/bar 0K

