# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ENOMEDIUM.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ENOMEDIUM

