# tests/auto/InvalidValue/2.31/good/Service_RootImage_--usr--sbin.target
[Unit]
Description=Foo
[Service]
RootImage=/usr/sbin

