# tests/auto/InvalidValue/2.31/good/Service_RootImage_--bin.target
[Unit]
Description=Foo
[Service]
RootImage=/bin

