# tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_28%C2%B5s.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=28µs

