# tests/auto/InvalidValue/2.31/good/Service_StartLimitBurst_0.target
[Unit]
Description=Foo
[Service]
StartLimitBurst=0

