# tests/auto/InvalidValue/2.31/good/Service_ProtectSystem_false.target
[Unit]
Description=Foo
[Service]
ProtectSystem=false

