# tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_22weeks.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=22weeks

