# tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_SIGURG.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=SIGURG

