# tests/auto/InvalidValue/2.31/good/Service_ExecStop_--bin--sh.target
[Unit]
Description=Foo
[Service]
ExecStop=/bin/sh

