# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ENXIO.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ENXIO

