# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGABRT.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGABRT

