# tests/auto/InvalidValue/2.31/good/Service_GuessMainPID_false.target
[Unit]
Description=Foo
[Service]
GuessMainPID=false

