# tests/auto/InvalidValue/2.31/good/Service_ProtectHome_false.target
[Unit]
Description=Foo
[Service]
ProtectHome=false

