# tests/auto/InvalidValue/2.31/good/Service_IOSchedulingClass_none.target
[Unit]
Description=Foo
[Service]
IOSchedulingClass=none

