# tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGPWR.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGPWR

