# tests/auto/InvalidValue/2.31/good/Service_StandardError_inherit.target
[Unit]
Description=Foo
[Service]
StandardError=inherit

