# tests/auto/InvalidValue/2.31/good/Service_MountFlags_slave.target
[Unit]
Description=Foo
[Service]
MountFlags=slave

