# tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_27.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=27

