# tests/auto/InvalidValue/2.31/good/Service_Restart_on-abort.target
[Unit]
Description=Foo
[Service]
Restart=on-abort

