# tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGCHLD.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGCHLD

