# tests/auto/InvalidValue/2.31/good/Service_TTYReset_on.target
[Unit]
Description=Foo
[Service]
TTYReset=on

