# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EALREADY.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EALREADY

