# tests/auto/InvalidValue/2.31/good/Service_StandardError_socket.target
[Unit]
Description=Foo
[Service]
StandardError=socket

