# tests/auto/InvalidValue/2.31/good/Service_LimitMSGQUEUE_infinity.target
[Unit]
Description=Foo
[Service]
LimitMSGQUEUE=infinity

