# tests/auto/InvalidValue/2.31/good/Service_SELinuxContext_abc.target
[Unit]
Description=Foo
[Service]
SELinuxContext=abc

