# tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_SIGSYS.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=SIGSYS

