# tests/auto/InvalidValue/2.31/good/Service_DevicePolicy_auto.target
[Unit]
Description=Foo
[Service]
DevicePolicy=auto

