# tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_SIGCLD.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=SIGCLD

