# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_15ms.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=15ms

