# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_9min.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=9min

