# tests/auto/InvalidValue/2.31/good/Service_PrivateDevices_yes.target
[Unit]
Description=Foo
[Service]
PrivateDevices=yes

