# tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_7minutes.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=7minutes

