# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGXFSZ.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGXFSZ

