# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ENOTCONN.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ENOTCONN

