# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGUSR1.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGUSR1

