# tests/auto/InvalidValue/2.31/good/Service_PrivateDevices_false.target
[Unit]
Description=Foo
[Service]
PrivateDevices=false

