# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGPOLL.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGPOLL

