# tests/auto/InvalidValue/2.31/good/Service_MountAPIVFS_false.target
[Unit]
Description=Foo
[Service]
MountAPIVFS=false

