# tests/auto/InvalidValue/2.31/good/Service_LimitCORE_9999999999999.target
[Unit]
Description=Foo
[Service]
LimitCORE=9999999999999

