# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGHUP.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGHUP

