# tests/auto/InvalidValue/2.31/good/Service_SyslogIdentifier_abc.target
[Unit]
Description=Foo
[Service]
SyslogIdentifier=abc

