# tests/auto/InvalidValue/2.31/good/Service_RestartSec_20usec.target
[Unit]
Description=Foo
[Service]
RestartSec=20usec

