# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_13months.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=13months

