# tests/auto/InvalidValue/2.31/good/Service_PIDFile_--bin.target
[Unit]
Description=Foo
[Service]
PIDFile=/bin

