# tests/auto/InvalidValue/2.31/good/Service_ProtectSystem_true.target
[Unit]
Description=Foo
[Service]
ProtectSystem=true

