# tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGPIPE.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGPIPE

