# tests/auto/InvalidValue/2.31/good/Service_NotifyAccess_all.target
[Unit]
Description=Foo
[Service]
NotifyAccess=all

