# tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGIO.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGIO

