# tests/auto/InvalidValue/2.31/good/Service_RemainAfterExit_1.target
[Unit]
Description=Foo
[Service]
RemainAfterExit=1

