# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_58.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=58

