# tests/auto/InvalidValue/2.31/good/Service_RemainAfterExit_off.target
[Unit]
Description=Foo
[Service]
RemainAfterExit=off

