# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ENOPKG.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ENOPKG

