# tests/auto/InvalidValue/2.31/good/Service_Personality_ppc-le.target
[Unit]
Description=Foo
[Service]
Personality=ppc-le

