# tests/auto/InvalidValue/2.31/good/Service_TTYVHangup_0.target
[Unit]
Description=Foo
[Service]
TTYVHangup=0

