# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EPERM.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EPERM

