# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGLOST.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGLOST

