# tests/auto/InvalidValue/2.31/good/Service_RestartSec_12month.target
[Unit]
Description=Foo
[Service]
RestartSec=12month

