# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ESTRPIPE.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ESTRPIPE

