# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EPROTOTYPE.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EPROTOTYPE

