# tests/auto/InvalidValue/2.31/good/Service_RestartSec_9min.target
[Unit]
Description=Foo
[Service]
RestartSec=9min

