# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EUNATCH.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EUNATCH

