# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_23week.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=23week

