# tests/auto/InvalidValue/2.31/good/Service_RemoveIPC_false.target
[Unit]
Description=Foo
[Service]
RemoveIPC=false

