# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_21.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=21

