# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ELIBSCN.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ELIBSCN

