# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_4hour.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=4hour

