# tests/auto/InvalidValue/2.31/good/Service_LimitDATA_0.target
[Unit]
Description=Foo
[Service]
LimitDATA=0

