# tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGPWR.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGPWR

