# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_18sec.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=18sec

