# tests/auto/InvalidValue/2.31/good/Service_IOSchedulingPriority_7.target
[Unit]
Description=Foo
[Service]
IOSchedulingPriority=7

