# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EBADF.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EBADF

