# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_DATAERR.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=DATAERR

