# tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGTTOU.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGTTOU

