# tests/auto/InvalidValue/2.31/good/Service_StartupCPUShares_2.target
[Unit]
Description=Foo
[Service]
StartupCPUShares=2

