# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGWINCH.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGWINCH

