# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGCLD.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGCLD

