# tests/auto/InvalidValue/2.31/good/Service_PrivateNetwork_on.target
[Unit]
Description=Foo
[Service]
PrivateNetwork=on

