# tests/auto/InvalidValue/2.31/good/Service_MemoryLow_0K.target
[Unit]
Description=Foo
[Service]
MemoryLow=0K

