# tests/auto/InvalidValue/2.31/good/Service_LimitRTTIME_9999999999999.target
[Unit]
Description=Foo
[Service]
LimitRTTIME=9999999999999

