# tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_SIGPROF.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=SIGPROF

