# tests/auto/InvalidValue/2.31/good/Service_StandardOutput_journal.target
[Unit]
Description=Foo
[Service]
StandardOutput=journal

