# tests/auto/InvalidValue/2.31/good/Service_MountAPIVFS_0.target
[Unit]
Description=Foo
[Service]
MountAPIVFS=0

