# tests/auto/InvalidValue/2.31/good/Service_RemoveIPC_1.target
[Unit]
Description=Foo
[Service]
RemoveIPC=1

