# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGUNUSED.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGUNUSED

