# tests/auto/InvalidValue/2.31/good/Service_CPUSchedulingResetOnFork_no.target
[Unit]
Description=Foo
[Service]
CPUSchedulingResetOnFork=no

