# tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGXFSZ.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGXFSZ

