# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_10.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=10

