# tests/auto/InvalidValue/2.31/good/Service_TTYReset_yes.target
[Unit]
Description=Foo
[Service]
TTYReset=yes

