# tests/auto/InvalidValue/2.31/good/Service_RestrictRealtime_no.target
[Unit]
Description=Foo
[Service]
RestrictRealtime=no

