# tests/auto/InvalidValue/2.31/good/Service_PassEnvironment_abc.target
[Unit]
Description=Foo
[Service]
PassEnvironment=abc

