# tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_6h.target
[Unit]
Description=Foo
[Service]
WatchdogSec=6h

