# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGINT.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGINT

