# tests/auto/InvalidValue/2.31/good/Service_IgnoreSIGPIPE_yes.target
[Unit]
Description=Foo
[Service]
IgnoreSIGPIPE=yes

