# tests/auto/InvalidValue/2.31/good/Service_LimitRTPRIO_9999999999999.target
[Unit]
Description=Foo
[Service]
LimitRTPRIO=9999999999999

