# tests/auto/InvalidValue/2.31/good/Service_LimitNPROC_9999999999999.target
[Unit]
Description=Foo
[Service]
LimitNPROC=9999999999999

