# tests/auto/InvalidValue/2.31/good/Service_BlockIODeviceWeight_--dev--foo%201.target
[Unit]
Description=Foo
[Service]
BlockIODeviceWeight=/dev/foo 1

