# tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_10M.target
[Unit]
Description=Foo
[Service]
TimeoutSec=10M

