# tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_13months.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=13months

