# tests/auto/InvalidValue/2.31/good/Service_ProtectHome_no.target
[Unit]
Description=Foo
[Service]
ProtectHome=no

