# tests/auto/InvalidValue/2.31/good/Service_CPUQuota_0%25.target
[Unit]
Description=Foo
[Service]
CPUQuota=0%

