# tests/auto/InvalidValue/2.31/good/Service_Restart_on-failure.target
[Unit]
Description=Foo
[Service]
Restart=on-failure

