# tests/auto/InvalidValue/2.31/good/Service_CapabilityBoundingSet_abc.target
[Unit]
Description=Foo
[Service]
CapabilityBoundingSet=abc

