# tests/auto/InvalidValue/2.31/good/Service_GuessMainPID_no.target
[Unit]
Description=Foo
[Service]
GuessMainPID=no

