# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EL2NSYNC.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EL2NSYNC

