# tests/auto/InvalidValue/2.31/good/Service_KillMode_none.target
[Unit]
Description=Foo
[Service]
KillMode=none

