# tests/auto/InvalidValue/2.31/good/Service_CPUSchedulingResetOnFork_off.target
[Unit]
Description=Foo
[Service]
CPUSchedulingResetOnFork=off

