# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGXFSZ.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGXFSZ

