# tests/auto/InvalidValue/2.31/good/Service_BlockIOAccounting_no.target
[Unit]
Description=Foo
[Service]
BlockIOAccounting=no

