# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGUNUSED.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGUNUSED

