# tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_17second.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=17second

