# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGQUIT.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGQUIT

