# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_infinity.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=infinity

