# tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_10M.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=10M

