# tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_22weeks.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=22weeks

