# tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_9min.target
[Unit]
Description=Foo
[Service]
TimeoutSec=9min

