# tests/auto/InvalidValue/2.31/good/Service_SyslogLevelPrefix_1.target
[Unit]
Description=Foo
[Service]
SyslogLevelPrefix=1

