# tests/auto/InvalidValue/2.31/good/Service_NotifyAccess_main.target
[Unit]
Description=Foo
[Service]
NotifyAccess=main

