# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ERFKILL.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ERFKILL

