# tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGEMT.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGEMT

