# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EREMOTEIO.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EREMOTEIO

