# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGPIPE.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGPIPE

