# tests/auto/InvalidValue/2.31/good/Service_LimitMEMLOCK_9999999999999.target
[Unit]
Description=Foo
[Service]
LimitMEMLOCK=9999999999999

