# tests/auto/InvalidValue/2.31/good/Service_Personality_x86-64.target
[Unit]
Description=Foo
[Service]
Personality=x86-64

