# tests/auto/InvalidValue/2.31/good/Service_CPUSchedulingResetOnFork_yes.target
[Unit]
Description=Foo
[Service]
CPUSchedulingResetOnFork=yes

