# tests/auto/InvalidValue/2.31/good/Service_UMask_9999999999999.target
[Unit]
Description=Foo
[Service]
UMask=9999999999999

