# tests/auto/InvalidValue/2.31/good/Service_StartLimitIntervalSec_22weeks.target
[Unit]
Description=Foo
[Service]
StartLimitIntervalSec=22weeks

