# tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_4hour.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=4hour

