# tests/auto/InvalidValue/2.31/good/Service_EnvironmentFile_--dev--.target
[Unit]
Description=Foo
[Service]
EnvironmentFile=/dev/

