# tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_DATAERR.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=DATAERR

