# tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_26year.target
[Unit]
Description=Foo
[Service]
TimeoutSec=26year

