# tests/auto/InvalidValue/2.31/good/Service_ProtectSystem_1.target
[Unit]
Description=Foo
[Service]
ProtectSystem=1

