# tests/auto/InvalidValue/2.31/good/Service_DisableControllers_cpu.target
[Unit]
Description=Foo
[Service]
DisableControllers=cpu

