# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_20usec.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=20usec

