# tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_18sec.target
[Unit]
Description=Foo
[Service]
WatchdogSec=18sec

