# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGPIPE.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGPIPE

