# tests/auto/InvalidValue/2.31/good/Service_TTYReset_0.target
[Unit]
Description=Foo
[Service]
TTYReset=0

