# tests/auto/InvalidValue/2.31/good/Service_TTYReset_off.target
[Unit]
Description=Foo
[Service]
TTYReset=off

