# tests/auto/InvalidValue/2.31/good/Service_RestrictRealtime_yes.target
[Unit]
Description=Foo
[Service]
RestrictRealtime=yes

