# tests/auto/InvalidValue/2.31/good/Service_LimitDATA_infinity.target
[Unit]
Description=Foo
[Service]
LimitDATA=infinity

