# tests/auto/InvalidValue/2.31/good/Service_StartLimitIntervalSec_13months.target
[Unit]
Description=Foo
[Service]
StartLimitIntervalSec=13months

