# tests/auto/InvalidValue/2.31/good/Service_StartLimitIntervalSec_21us.target
[Unit]
Description=Foo
[Service]
StartLimitIntervalSec=21us

