# tests/auto/InvalidValue/2.31/good/Service_Delegate_blkio.target
[Unit]
Description=Foo
[Service]
Delegate=blkio

