# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_49.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=49

