# tests/auto/InvalidValue/2.31/good/Service_MountAPIVFS_true.target
[Unit]
Description=Foo
[Service]
MountAPIVFS=true

