# tests/auto/InvalidValue/2.31/good/Service_ProtectSystem_0.target
[Unit]
Description=Foo
[Service]
ProtectSystem=0

