# tests/auto/InvalidValue/2.31/good/Service_IODeviceLatencyTargetSec_abc.target
[Unit]
Description=Foo
[Service]
IODeviceLatencyTargetSec=abc

