# tests/auto/InvalidValue/2.31/good/Service_ReadOnlyDirectories_abc.target
[Unit]
Description=Foo
[Service]
ReadOnlyDirectories=abc

