# tests/auto/InvalidValue/2.31/good/Service_SyslogLevelPrefix_false.target
[Unit]
Description=Foo
[Service]
SyslogLevelPrefix=false

