# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_14.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=14

