# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_4095.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=4095

