# tests/auto/InvalidValue/2.31/good/Service_Type_forking.target
[Unit]
Description=Foo
[Service]
Type=forking

