# tests/auto/InvalidValue/2.31/good/Service_CPUQuota_9999999999999%25.target
[Unit]
Description=Foo
[Service]
CPUQuota=9999999999999%

