# tests/auto/InvalidValue/2.31/good/Service_IOAccounting_no.target
[Unit]
Description=Foo
[Service]
IOAccounting=no

