# tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_24.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=24

