# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EACCES.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EACCES

