# tests/auto/InvalidValue/2.31/good/Service_MemoryMax_0G.target
[Unit]
Description=Foo
[Service]
MemoryMax=0G

