# tests/auto/InvalidValue/2.31/good/Service_ProtectHome_off.target
[Unit]
Description=Foo
[Service]
ProtectHome=off

