# tests/auto/InvalidValue/2.31/good/Service_MemoryAccounting_on.target
[Unit]
Description=Foo
[Service]
MemoryAccounting=on

