# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_24w.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=24w

