# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ELIBBAD.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ELIBBAD

