# tests/auto/InvalidValue/2.31/good/Service_StandardOutput_socket.target
[Unit]
Description=Foo
[Service]
StandardOutput=socket

