# tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGSEGV.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGSEGV

