# tests/auto/InvalidValue/2.31/good/Service_LimitNPROC_infinity.target
[Unit]
Description=Foo
[Service]
LimitNPROC=infinity

