# tests/auto/InvalidValue/2.31/good/Service_SecureBits_keep-caps-locked.target
[Unit]
Description=Foo
[Service]
SecureBits=keep-caps-locked

