# tests/auto/InvalidValue/2.31/good/Service_RestrictSUIDSGID_0.target
[Unit]
Description=Foo
[Service]
RestrictSUIDSGID=0

