# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EHWPOISON.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EHWPOISON

