# tests/auto/InvalidValue/2.31/good/Service_PrivateNetwork_no.target
[Unit]
Description=Foo
[Service]
PrivateNetwork=no

