# tests/auto/InvalidValue/2.31/good/Service_MemorySwapMax_0G.target
[Unit]
Description=Foo
[Service]
MemorySwapMax=0G

