# tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_21us.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=21us

