# tests/auto/InvalidValue/2.31/good/Service_TasksAccounting_true.target
[Unit]
Description=Foo
[Service]
TasksAccounting=true

