# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_1.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=1

