# tests/auto/InvalidValue/2.31/good/Service_BlockIOWriteBandwidthMax_--dev--bar%200G.target
[Unit]
Description=Foo
[Service]
BlockIOWriteBandwidthMax=/dev/bar 0G

