# tests/auto/InvalidValue/2.31/good/Service_BlockIOAccounting_1.target
[Unit]
Description=Foo
[Service]
BlockIOAccounting=1

