# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_5hr.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=5hr

