# tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGURG.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGURG

