# tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_16.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=16

