# tests/auto/InvalidValue/2.31/good/Service_IOWriteBandwidthMax_abc.target
[Unit]
Description=Foo
[Service]
IOWriteBandwidthMax=abc

