# tests/auto/InvalidValue/2.31/good/Service_StandardOutput_kmsg%2Bconsole.target
[Unit]
Description=Foo
[Service]
StandardOutput=kmsg+console

