# tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGPROF.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGPROF

