# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGVTALRM.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGVTALRM

