# tests/auto/InvalidValue/2.31/good/Service_GuessMainPID_off.target
[Unit]
Description=Foo
[Service]
GuessMainPID=off

