# tests/auto/InvalidValue/2.31/good/Service_MemoryAccounting_yes.target
[Unit]
Description=Foo
[Service]
MemoryAccounting=yes

