# tests/auto/InvalidValue/2.31/good/Service_Delegate_devices.target
[Unit]
Description=Foo
[Service]
Delegate=devices

