# tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_SIGTSTP.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=SIGTSTP

