# tests/auto/InvalidValue/2.31/good/Service_RemoveIPC_true.target
[Unit]
Description=Foo
[Service]
RemoveIPC=true

