# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_27y.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=27y

