# tests/auto/InvalidValue/2.31/good/Service_TasksAccounting_0.target
[Unit]
Description=Foo
[Service]
TasksAccounting=0

