# tests/auto/InvalidValue/2.31/good/Service_RemoveIPC_no.target
[Unit]
Description=Foo
[Service]
RemoveIPC=no

