# tests/auto/InvalidValue/2.31/good/Service_SyslogLevelPrefix_on.target
[Unit]
Description=Foo
[Service]
SyslogLevelPrefix=on

