# tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_16seconds.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=16seconds

