# tests/auto/InvalidValue/2.31/good/Service_LimitNICE_0.target
[Unit]
Description=Foo
[Service]
LimitNICE=0

