# tests/auto/InvalidValue/2.31/good/Service_UtmpMode_user.target
[Unit]
Description=Foo
[Service]
UtmpMode=user

