# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGINT.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGINT

