# tests/auto/InvalidValue/2.31/good/Service_TTYVTDisallocate_on.target
[Unit]
Description=Foo
[Service]
TTYVTDisallocate=on

