# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_17second.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=17second

