# tests/auto/InvalidValue/2.31/good/Service_LimitRTPRIO_infinity.target
[Unit]
Description=Foo
[Service]
LimitRTPRIO=infinity

