# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGURG.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGURG

