# tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGSTOP.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGSTOP

