# tests/auto/InvalidValue/2.31/good/Service_User_metoo.target
[Unit]
Description=Foo
[Service]
User=metoo

