# tests/auto/InvalidValue/2.31/good/Service_MemoryMax_9999999999999M.target
[Unit]
Description=Foo
[Service]
MemoryMax=9999999999999M

