# tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_5hr.target
[Unit]
Description=Foo
[Service]
TimeoutSec=5hr

