# tests/auto/InvalidValue/2.31/good/Service_AppArmorProfile_abc.target
[Unit]
Description=Foo
[Service]
AppArmorProfile=abc

