# tests/auto/InvalidValue/2.31/good/Service_Type_exec.target
[Unit]
Description=Foo
[Service]
Type=exec

