# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ENOSR.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ENOSR

