# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_31.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=31

