# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGIO.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGIO

