# tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_11m.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=11m

