# tests/auto/InvalidValue/2.31/good/Service_SendSIGKILL_off.target
[Unit]
Description=Foo
[Service]
SendSIGKILL=off

