# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGSYS.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGSYS

