# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ENODATA.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ENODATA

