# tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_9min.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=9min

