# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGIOT.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGIOT

