# tests/auto/InvalidValue/2.31/good/Service_TasksAccounting_false.target
[Unit]
Description=Foo
[Service]
TasksAccounting=false

