# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EBADSLT.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EBADSLT

