# tests/auto/InvalidValue/2.31/good/Service_StandardError_kmsg%2Bconsole.target
[Unit]
Description=Foo
[Service]
StandardError=kmsg+console

