# tests/auto/InvalidValue/2.31/good/Service_SendSIGHUP_false.target
[Unit]
Description=Foo
[Service]
SendSIGHUP=false

