# tests/auto/InvalidValue/2.31/good/Service_SyslogFacility_local4.target
[Unit]
Description=Foo
[Service]
SyslogFacility=local4

