# tests/auto/InvalidValue/2.31/good/Service_SecureBits_noroot.target
[Unit]
Description=Foo
[Service]
SecureBits=noroot

