# tests/auto/InvalidValue/2.31/good/Service_SyslogLevel_warning.target
[Unit]
Description=Foo
[Service]
SyslogLevel=warning

