# tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_14msec.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=14msec

