# tests/auto/InvalidValue/2.31/good/Service_IOAccounting_yes.target
[Unit]
Description=Foo
[Service]
IOAccounting=yes

