# tests/auto/InvalidValue/2.31/good/Service_RestartSec_2d.target
[Unit]
Description=Foo
[Service]
RestartSec=2d

