# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_16seconds.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=16seconds

