# tests/auto/InvalidValue/2.31/good/Service_TTYVHangup_no.target
[Unit]
Description=Foo
[Service]
TTYVHangup=no

