# tests/auto/InvalidValue/2.31/good/Service_MemorySwapMax_9999999999999T.target
[Unit]
Description=Foo
[Service]
MemorySwapMax=9999999999999T

