# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EEXIST.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EEXIST

