# tests/auto/InvalidValue/2.31/good/Service_SyslogFacility_lpr.target
[Unit]
Description=Foo
[Service]
SyslogFacility=lpr

