# tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_SIGXFSZ.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=SIGXFSZ

