# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGTTOU.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGTTOU

