# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGTSTP.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGTSTP

