# tests/auto/InvalidValue/2.31/good/Service_SecureBits_no-setuid-fixup.target
[Unit]
Description=Foo
[Service]
SecureBits=no-setuid-fixup

