# tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_17second.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=17second

