# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ERANGE.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ERANGE

