# tests/auto/InvalidValue/2.31/good/Service_USBFunctionStrings_abc.target
[Unit]
Description=Foo
[Service]
USBFunctionStrings=abc

