# tests/auto/InvalidValue/2.31/good/Service_StandardOutput_null.target
[Unit]
Description=Foo
[Service]
StandardOutput=null

