# tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_14msec.target
[Unit]
Description=Foo
[Service]
WatchdogSec=14msec

