# tests/auto/InvalidValue/2.31/good/Service_DeviceAllow_--dev--foo.target
[Unit]
Description=Foo
[Service]
DeviceAllow=/dev/foo

