# tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_18sec.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=18sec

