# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_42.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=42

