# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_19.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=19

