# tests/auto/InvalidValue/2.31/good/Service_LimitMEMLOCK_0.target
[Unit]
Description=Foo
[Service]
LimitMEMLOCK=0

