# tests/auto/InvalidValue/2.31/good/Service_RestartSec_infinity.target
[Unit]
Description=Foo
[Service]
RestartSec=infinity

