# tests/auto/InvalidValue/2.31/good/Service_SendSIGKILL_0.target
[Unit]
Description=Foo
[Service]
SendSIGKILL=0

