# tests/auto/InvalidValue/2.31/good/Service_LimitRTTIME_infinity.target
[Unit]
Description=Foo
[Service]
LimitRTTIME=infinity

