# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_52.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=52

