# tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGTERM.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGTERM

