# tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_9min.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=9min

