# tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_6h.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=6h

