# tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGPOLL.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGPOLL

