# tests/auto/InvalidValue/2.31/good/Service_BlockIOWriteBandwidthMax_--dev--foo%200K.target
[Unit]
Description=Foo
[Service]
BlockIOWriteBandwidthMax=/dev/foo 0K

