# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGEMT.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGEMT

