# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EADDRNOTAVAIL.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EADDRNOTAVAIL

