# tests/auto/InvalidValue/2.31/good/Service_PrivateNetwork_off.target
[Unit]
Description=Foo
[Service]
PrivateNetwork=off

