# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_E2BIG.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=E2BIG

