# tests/auto/InvalidValue/2.31/good/Service_Type_dbus.target
[Unit]
Description=Foo
[Service]
Type=dbus

