# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_2.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=2

