# tests/auto/InvalidValue/2.31/good/Service_RestartSec_5hr.target
[Unit]
Description=Foo
[Service]
RestartSec=5hr

