# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EDOM.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EDOM

