# tests/auto/InvalidValue/2.31/good/Service_CPUSchedulingPriority_1.target
[Unit]
Description=Foo
[Service]
CPUSchedulingPriority=1

