# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGUSR1.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGUSR1

