# tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGFPE.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGFPE

