# tests/auto/InvalidValue/2.31/good/Service_PrivateTmp_0.target
[Unit]
Description=Foo
[Service]
PrivateTmp=0

