# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGCHLD.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGCHLD

