# tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_26year.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=26year

