# tests/auto/InvalidValue/2.31/good/Service_Type_notify.target
[Unit]
Description=Foo
[Service]
Type=notify

