# tests/auto/InvalidValue/2.31/good/Service_TTYVTDisallocate_true.target
[Unit]
Description=Foo
[Service]
TTYVTDisallocate=true

