# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_37.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=37

