# tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_19s.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=19s

