# tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_SIGUSR1.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=SIGUSR1

