# tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_42.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=42

