# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_38.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=38

