# tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_41.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=41

