# tests/auto/InvalidValue/2.31/good/Service_MemoryAccounting_1.target
[Unit]
Description=Foo
[Service]
MemoryAccounting=1

