# tests/auto/InvalidValue/2.31/good/Service_DevicePolicy_closed.target
[Unit]
Description=Foo
[Service]
DevicePolicy=closed

