# tests/auto/InvalidValue/2.31/good/Service_IgnoreSIGPIPE_no.target
[Unit]
Description=Foo
[Service]
IgnoreSIGPIPE=no

