# tests/auto/InvalidValue/2.31/good/Service_RootDirectory_--usr--sbin.target
[Unit]
Description=Foo
[Service]
RootDirectory=/usr/sbin

