# tests/auto/InvalidValue/2.31/good/Service_KillSignal_SIGXCPU.target
[Unit]
Description=Foo
[Service]
KillSignal=SIGXCPU

