# tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_20usec.target
[Unit]
Description=Foo
[Service]
TimeoutSec=20usec

