# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGQUIT.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGQUIT

