# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGSTKFLT.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGSTKFLT

