# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_1day.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=1day

