# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGLOST.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGLOST

