# tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGVTALRM.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGVTALRM

