# tests/auto/InvalidValue/2.31/good/Service_LimitLOCKS_infinity.target
[Unit]
Description=Foo
[Service]
LimitLOCKS=infinity

