# tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_25years.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=25years

