# tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_23week.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=23week

