# tests/auto/InvalidValue/2.31/good/Service_ProtectSystem_off.target
[Unit]
Description=Foo
[Service]
ProtectSystem=off

