# tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGTRAP.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGTRAP

