# tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_13months.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=13months

