# tests/auto/InvalidValue/2.31/good/Service_TimeoutStartSec_25years.target
[Unit]
Description=Foo
[Service]
TimeoutStartSec=25years

