# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ERESTART.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ERESTART

