# tests/auto/InvalidValue/2.31/good/Service_SendSIGHUP_off.target
[Unit]
Description=Foo
[Service]
SendSIGHUP=off

