# tests/auto/InvalidValue/2.31/good/Service_Delegate_memory.target
[Unit]
Description=Foo
[Service]
Delegate=memory

