# tests/auto/InvalidValue/2.31/good/Service_NonBlocking_0.target
[Unit]
Description=Foo
[Service]
NonBlocking=0

