# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EAGAIN.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EAGAIN

