# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ECHILD.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ECHILD

