# tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_14.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=14

