# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGLOST.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGLOST

