# tests/auto/InvalidValue/2.31/good/Service_SyslogLevel_alert.target
[Unit]
Description=Foo
[Service]
SyslogLevel=alert

