# tests/auto/InvalidValue/2.31/good/Service_DevicePolicy_strict.target
[Unit]
Description=Foo
[Service]
DevicePolicy=strict

