# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGCHLD.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGCHLD

