# tests/auto/InvalidValue/2.31/good/Service_LimitCORE_infinity.target
[Unit]
Description=Foo
[Service]
LimitCORE=infinity

