# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EKEYEXPIRED.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EKEYEXPIRED

