# tests/auto/InvalidValue/2.31/good/Service_StartLimitIntervalSec_11m.target
[Unit]
Description=Foo
[Service]
StartLimitIntervalSec=11m

