# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGPWR.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGPWR

