# tests/auto/InvalidValue/2.31/good/Service_UtmpMode_login.target
[Unit]
Description=Foo
[Service]
UtmpMode=login

