# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EHOSTDOWN.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EHOSTDOWN

