# tests/auto/InvalidValue/2.31/good/Service_USBFunctionDescriptors_abc.target
[Unit]
Description=Foo
[Service]
USBFunctionDescriptors=abc

