# tests/auto/InvalidValue/2.31/good/Service_LimitDATA_9999999999999.target
[Unit]
Description=Foo
[Service]
LimitDATA=9999999999999

