# tests/auto/InvalidValue/2.31/good/Service_TTYVTDisallocate_0.target
[Unit]
Description=Foo
[Service]
TTYVTDisallocate=0

