# tests/auto/InvalidValue/2.31/good/Service_SyslogLevelPrefix_no.target
[Unit]
Description=Foo
[Service]
SyslogLevelPrefix=no

