# tests/auto/InvalidValue/2.31/good/Service_KillSignal_13.target
[Unit]
Description=Foo
[Service]
KillSignal=13

