# tests/auto/InvalidValue/2.31/good/Service_RuntimeMaxSec_10M.target
[Unit]
Description=Foo
[Service]
RuntimeMaxSec=10M

