# tests/auto/InvalidValue/2.31/good/Service_IOSchedulingClass_idle.target
[Unit]
Description=Foo
[Service]
IOSchedulingClass=idle

