# tests/auto/InvalidValue/2.31/good/Service_RemainAfterExit_on.target
[Unit]
Description=Foo
[Service]
RemainAfterExit=on

