# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGCLD.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGCLD

