# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_10M.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=10M

