# tests/auto/InvalidValue/2.31/good/Service_StartLimitIntervalSec_9min.target
[Unit]
Description=Foo
[Service]
StartLimitIntervalSec=9min

