# tests/auto/InvalidValue/2.31/good/Service_IOReadIOPSMax_abc.target
[Unit]
Description=Foo
[Service]
IOReadIOPSMax=abc

