# tests/auto/InvalidValue/2.31/good/Service_SyslogFacility_daemon.target
[Unit]
Description=Foo
[Service]
SyslogFacility=daemon

