# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGALRM.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGALRM

