# tests/auto/InvalidValue/2.31/good/Service_Personality_ppc.target
[Unit]
Description=Foo
[Service]
Personality=ppc

