# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGBUS.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGBUS

