# tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_SIGIOT.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=SIGIOT

