# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EFAULT.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EFAULT

