# tests/auto/InvalidValue/2.31/good/Service_OOMScoreAdjust_1000.target
[Unit]
Description=Foo
[Service]
OOMScoreAdjust=1000

