# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EDEADLK.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EDEADLK

