# tests/auto/InvalidValue/2.31/good/Service_MemoryHigh_0M.target
[Unit]
Description=Foo
[Service]
MemoryHigh=0M

