# tests/auto/InvalidValue/2.31/good/Service_Personality_ppc64-le.target
[Unit]
Description=Foo
[Service]
Personality=ppc64-le

