# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGUSR2.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGUSR2

