# tests/auto/InvalidValue/2.31/good/Service_KillMode_process.target
[Unit]
Description=Foo
[Service]
KillMode=process

