# tests/auto/InvalidValue/2.31/good/Service_TimeoutSec_25years.target
[Unit]
Description=Foo
[Service]
TimeoutSec=25years

