# tests/auto/InvalidValue/2.31/good/Service_MountFlags_private.target
[Unit]
Description=Foo
[Service]
MountFlags=private

