# tests/auto/InvalidValue/2.31/good/Service_MemoryAccounting_true.target
[Unit]
Description=Foo
[Service]
MemoryAccounting=true

