# tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_SIGPWR.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=SIGPWR

