# tests/auto/InvalidValue/2.31/good/Service_MountAPIVFS_no.target
[Unit]
Description=Foo
[Service]
MountAPIVFS=no

