# tests/auto/InvalidValue/2.31/good/Service_LimitRSS_infinity.target
[Unit]
Description=Foo
[Service]
LimitRSS=infinity

