# tests/auto/InvalidValue/2.31/good/Service_WatchdogSignal_SIGHUP.target
[Unit]
Description=Foo
[Service]
WatchdogSignal=SIGHUP

