# tests/auto/InvalidValue/2.31/good/Service_PrivateNetwork_0.target
[Unit]
Description=Foo
[Service]
PrivateNetwork=0

