# tests/auto/InvalidValue/2.31/good/Service_Type_oneshot.target
[Unit]
Description=Foo
[Service]
Type=oneshot

