# tests/auto/InvalidValue/2.31/good/Service_SendSIGKILL_1.target
[Unit]
Description=Foo
[Service]
SendSIGKILL=1

