# tests/auto/InvalidValue/2.31/good/Service_Delegate_io.target
[Unit]
Description=Foo
[Service]
Delegate=io

