# tests/auto/InvalidValue/2.31/good/Service_SuccessExitStatus_18.target
[Unit]
Description=Foo
[Service]
SuccessExitStatus=18

