# tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_28%C2%B5s.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=28µs

