# tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_infinity.target
[Unit]
Description=Foo
[Service]
WatchdogSec=infinity

