# tests/auto/InvalidValue/2.31/good/Service_RestartPreventExitStatus_SIGFPE.target
[Unit]
Description=Foo
[Service]
RestartPreventExitStatus=SIGFPE

