# tests/auto/InvalidValue/2.31/good/Service_TTYVHangup_false.target
[Unit]
Description=Foo
[Service]
TTYVHangup=false

