# tests/auto/InvalidValue/2.31/good/Service_TimerAccuracySec_25years.target
[Unit]
Description=Foo
[Service]
TimerAccuracySec=25years

