# tests/auto/InvalidValue/2.31/good/Service_MountAPIVFS_off.target
[Unit]
Description=Foo
[Service]
MountAPIVFS=off

