# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGILL.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGILL

