# tests/auto/InvalidValue/2.31/good/Service_TimerSlackNSec_8minute.target
[Unit]
Description=Foo
[Service]
TimerSlackNSec=8minute

