# tests/auto/InvalidValue/2.31/good/Service_FinalKillSignal_16.target
[Unit]
Description=Foo
[Service]
FinalKillSignal=16

