# tests/auto/InvalidValue/2.31/good/Service_CPUSchedulingPolicy_batch.target
[Unit]
Description=Foo
[Service]
CPUSchedulingPolicy=batch

