# tests/auto/InvalidValue/2.31/good/Service_TimeoutStopSec_20usec.target
[Unit]
Description=Foo
[Service]
TimeoutStopSec=20usec

