# tests/auto/InvalidValue/2.31/good/Service_RestartForceExitStatus_SIGSEGV.target
[Unit]
Description=Foo
[Service]
RestartForceExitStatus=SIGSEGV

