# tests/auto/InvalidValue/2.31/good/Service_Restart_no.target
[Unit]
Description=Foo
[Service]
Restart=no

