# tests/auto/InvalidValue/2.31/good/Service_BlockIOWriteBandwidthMax_--dev--foo%209999999999999M.target
[Unit]
Description=Foo
[Service]
BlockIOWriteBandwidthMax=/dev/foo 9999999999999M

