# tests/auto/InvalidValue/2.31/good/Service_GuessMainPID_on.target
[Unit]
Description=Foo
[Service]
GuessMainPID=on

