# tests/auto/InvalidValue/2.31/good/Service_StartLimitIntervalSec_16seconds.target
[Unit]
Description=Foo
[Service]
StartLimitIntervalSec=16seconds

