# tests/auto/InvalidValue/2.31/good/Service_RestrictRealtime_1.target
[Unit]
Description=Foo
[Service]
RestrictRealtime=1

