# tests/auto/InvalidValue/2.31/good/Service_ProtectHome_0.target
[Unit]
Description=Foo
[Service]
ProtectHome=0

