# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ENOTRECOVERABLE.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ENOTRECOVERABLE

