# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_ELIBACC.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=ELIBACC

