# tests/auto/InvalidValue/2.31/good/Service_SendSIGKILL_no.target
[Unit]
Description=Foo
[Service]
SendSIGKILL=no

