# tests/auto/InvalidValue/2.31/good/Service_TasksAccounting_off.target
[Unit]
Description=Foo
[Service]
TasksAccounting=off

