# tests/auto/InvalidValue/2.31/good/Service_ProtectSystem_full.target
[Unit]
Description=Foo
[Service]
ProtectSystem=full

