# tests/auto/InvalidValue/2.31/good/Service_SystemCallErrorNumber_EOVERFLOW.target
[Unit]
Description=Foo
[Service]
SystemCallErrorNumber=EOVERFLOW

