# tests/auto/InvalidValue/2.31/good/Service_WatchdogSec_15ms.target
[Unit]
Description=Foo
[Service]
WatchdogSec=15ms

