# tests/auto/InvalidValue/2.31/good/Socket_SELinuxContextFromNet_0.target
[Unit]
Description=Foo
[Socket]
SELinuxContextFromNet=0

