# tests/auto/InvalidValue/2.31/good/Socket_Delegate_cpu.target
[Unit]
Description=Foo
[Socket]
Delegate=cpu

