# tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_19s.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=19s

