# tests/auto/InvalidValue/2.31/good/Socket_SystemCallFilter_abc.target
[Unit]
Description=Foo
[Socket]
SystemCallFilter=abc

