# tests/auto/InvalidValue/2.31/good/Socket_Transparent_0.target
[Unit]
Description=Foo
[Socket]
Transparent=0

