# tests/auto/InvalidValue/2.31/good/Socket_SyslogLevel_emerg.target
[Unit]
Description=Foo
[Socket]
SyslogLevel=emerg

