# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_7.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=7

