# tests/auto/InvalidValue/2.31/good/Socket_IPTOS_throughput.target
[Unit]
Description=Foo
[Socket]
IPTOS=throughput

