# tests/auto/InvalidValue/2.31/good/Socket_ProtectSystem_1.target
[Unit]
Description=Foo
[Socket]
ProtectSystem=1

