# tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_15ms.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=15ms

