# tests/auto/InvalidValue/2.31/good/Socket_IgnoreSIGPIPE_true.target
[Unit]
Description=Foo
[Socket]
IgnoreSIGPIPE=true

