# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EDQUOT.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EDQUOT

