# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EFBIG.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EFBIG

