# tests/auto/InvalidValue/2.31/good/Socket_LimitNPROC_infinity.target
[Unit]
Description=Foo
[Socket]
LimitNPROC=infinity

