# tests/auto/InvalidValue/2.31/good/Socket_SyslogFacility_local5.target
[Unit]
Description=Foo
[Socket]
SyslogFacility=local5

