# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGINFO.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGINFO

