# tests/auto/InvalidValue/2.31/good/Socket_PassCredentials_0.target
[Unit]
Description=Foo
[Socket]
PassCredentials=0

