# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ELIBSCN.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ELIBSCN

