# tests/auto/InvalidValue/2.31/good/Socket_DeferAcceptSec_0.target
[Unit]
Description=Foo
[Socket]
DeferAcceptSec=0

