# tests/auto/InvalidValue/2.31/good/Socket_ProtectSystem_no.target
[Unit]
Description=Foo
[Socket]
ProtectSystem=no

