# tests/auto/InvalidValue/2.31/good/Socket_BlockIOWriteBandwidthMax_--dev--foo%200T.target
[Unit]
Description=Foo
[Socket]
BlockIOWriteBandwidthMax=/dev/foo 0T

