# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EBADF.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EBADF

