# tests/auto/InvalidValue/2.31/good/Socket_CPUAffinity_abc.target
[Unit]
Description=Foo
[Socket]
CPUAffinity=abc

