# tests/auto/InvalidValue/2.31/good/Socket_MemoryLow_0B.target
[Unit]
Description=Foo
[Socket]
MemoryLow=0B

