# tests/auto/InvalidValue/2.31/good/Socket_PrivateDevices_true.target
[Unit]
Description=Foo
[Socket]
PrivateDevices=true

