# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGTERM.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGTERM

