# tests/auto/InvalidValue/2.31/good/Socket_PrivateNetwork_off.target
[Unit]
Description=Foo
[Socket]
PrivateNetwork=off

