# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGUSR1.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGUSR1

