# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGFPE.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGFPE

