# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGPIPE.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGPIPE

