# tests/auto/InvalidValue/2.31/good/Socket_Writable_1.target
[Unit]
Description=Foo
[Socket]
Writable=1

