# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_40.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=40

