# tests/auto/InvalidValue/2.31/good/Socket_LimitNOFILE_0.target
[Unit]
Description=Foo
[Socket]
LimitNOFILE=0

