# tests/auto/InvalidValue/2.31/good/Socket_User_iamauser.target
[Unit]
Description=Foo
[Socket]
User=iamauser

