# tests/auto/InvalidValue/2.31/good/Socket_PrivateNetwork_yes.target
[Unit]
Description=Foo
[Socket]
PrivateNetwork=yes

