# tests/auto/InvalidValue/2.31/good/Socket_MemoryAccounting_true.target
[Unit]
Description=Foo
[Socket]
MemoryAccounting=true

