# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGINT.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGINT

