# tests/auto/InvalidValue/2.31/good/Socket_NoDelay_0.target
[Unit]
Description=Foo
[Socket]
NoDelay=0

