# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGCONT.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGCONT

