# tests/auto/InvalidValue/2.31/good/Socket_Broadcast_true.target
[Unit]
Description=Foo
[Socket]
Broadcast=true

