# tests/auto/InvalidValue/2.31/good/Socket_BindToDevice_--dev--bar.target
[Unit]
Description=Foo
[Socket]
BindToDevice=/dev/bar

