# tests/auto/InvalidValue/2.31/good/Socket_CPUSchedulingPolicy_other.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingPolicy=other

