# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGPROF.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGPROF

