# tests/auto/InvalidValue/2.31/good/Socket_SyslogFacility_syslog.target
[Unit]
Description=Foo
[Socket]
SyslogFacility=syslog

