# tests/auto/InvalidValue/2.31/good/Socket_ProtectHome_no.target
[Unit]
Description=Foo
[Socket]
ProtectHome=no

