# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGSEGV.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGSEGV

