# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EKEYREVOKED.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EKEYREVOKED

