# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENFILE.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENFILE

