# tests/auto/InvalidValue/2.31/good/Socket_PrivateNetwork_no.target
[Unit]
Description=Foo
[Socket]
PrivateNetwork=no

