# tests/auto/InvalidValue/2.31/good/Socket_PAMName_abc.target
[Unit]
Description=Foo
[Socket]
PAMName=abc

