# tests/auto/InvalidValue/2.31/good/Socket_ExecStartPre_--bin--sh.target
[Unit]
Description=Foo
[Socket]
ExecStartPre=/bin/sh

