# tests/auto/InvalidValue/2.31/good/Socket_TTYReset_1.target
[Unit]
Description=Foo
[Socket]
TTYReset=1

