# tests/auto/InvalidValue/2.31/good/Socket_TasksAccounting_off.target
[Unit]
Description=Foo
[Socket]
TasksAccounting=off

