# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_58.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=58

