# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EREMOTEIO.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EREMOTEIO

