# tests/auto/InvalidValue/2.31/good/Socket_LimitRSS_9999999999999.target
[Unit]
Description=Foo
[Socket]
LimitRSS=9999999999999

