# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ESTRPIPE.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ESTRPIPE

