# tests/auto/InvalidValue/2.31/good/Socket_LimitRSS_0.target
[Unit]
Description=Foo
[Socket]
LimitRSS=0

