# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ERFKILL.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ERFKILL

