# tests/auto/InvalidValue/2.31/good/Socket_CPUQuota_0%25.target
[Unit]
Description=Foo
[Socket]
CPUQuota=0%

