# tests/auto/InvalidValue/2.31/good/Socket_NoDelay_no.target
[Unit]
Description=Foo
[Socket]
NoDelay=no

