# tests/auto/InvalidValue/2.31/good/Socket_DevicePolicy_closed.target
[Unit]
Description=Foo
[Socket]
DevicePolicy=closed

