# tests/auto/InvalidValue/2.31/good/Socket_SyslogIdentifier_abc.target
[Unit]
Description=Foo
[Socket]
SyslogIdentifier=abc

