# tests/auto/InvalidValue/2.31/good/Socket_MemoryDenyWriteExecute_no.target
[Unit]
Description=Foo
[Socket]
MemoryDenyWriteExecute=no

