# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ECHILD.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ECHILD

