# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGCHLD.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGCHLD

