# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EPROTO.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EPROTO

