# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGABRT.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGABRT

