# tests/auto/InvalidValue/2.31/good/Socket_TTYReset_off.target
[Unit]
Description=Foo
[Socket]
TTYReset=off

