# tests/auto/InvalidValue/2.31/good/Socket_IOWriteBandwidthMax_abc.target
[Unit]
Description=Foo
[Socket]
IOWriteBandwidthMax=abc

