# tests/auto/InvalidValue/2.31/good/Socket_BlockIOAccounting_false.target
[Unit]
Description=Foo
[Socket]
BlockIOAccounting=false

