# tests/auto/InvalidValue/2.31/good/Socket_Broadcast_on.target
[Unit]
Description=Foo
[Socket]
Broadcast=on

