# tests/auto/InvalidValue/2.31/good/Socket_RemoveIPC_no.target
[Unit]
Description=Foo
[Socket]
RemoveIPC=no

