# tests/auto/InvalidValue/2.31/good/Socket_Nice_19.target
[Unit]
Description=Foo
[Socket]
Nice=19

