# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGINFO.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGINFO

