# tests/auto/InvalidValue/2.31/good/Socket_IOWeight_10000.target
[Unit]
Description=Foo
[Socket]
IOWeight=10000

