# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOTRECOVERABLE.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOTRECOVERABLE

