# tests/auto/InvalidValue/2.31/good/Socket_ReusePort_on.target
[Unit]
Description=Foo
[Socket]
ReusePort=on

