# tests/auto/InvalidValue/2.31/good/Socket_IOWriteIOPSMax_abc.target
[Unit]
Description=Foo
[Socket]
IOWriteIOPSMax=abc

