# tests/auto/InvalidValue/2.31/good/Socket_LimitLOCKS_infinity.target
[Unit]
Description=Foo
[Socket]
LimitLOCKS=infinity

