# tests/auto/InvalidValue/2.31/good/Socket_WorkingDirectory_--usr--sbin.target
[Unit]
Description=Foo
[Socket]
WorkingDirectory=/usr/sbin

