# tests/auto/InvalidValue/2.31/good/Socket_SendSIGHUP_false.target
[Unit]
Description=Foo
[Socket]
SendSIGHUP=false

