# tests/auto/InvalidValue/2.31/good/Socket_RootDirectory_--bin.target
[Unit]
Description=Foo
[Socket]
RootDirectory=/bin

