# tests/auto/InvalidValue/2.31/good/Socket_TCPCongestion_abc.target
[Unit]
Description=Foo
[Socket]
TCPCongestion=abc

