# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EIDRM.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EIDRM

