# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGSTOP.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGSTOP

