# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGCLD.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGCLD

