# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGCONT.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGCONT

