# tests/auto/InvalidValue/2.31/good/Socket_IOReadBandwidthMax_abc.target
[Unit]
Description=Foo
[Socket]
IOReadBandwidthMax=abc

