# tests/auto/InvalidValue/2.31/good/Socket_RemoveOnStop_off.target
[Unit]
Description=Foo
[Socket]
RemoveOnStop=off

