# tests/auto/InvalidValue/2.31/good/Socket_Delegate_pids.target
[Unit]
Description=Foo
[Socket]
Delegate=pids

