# tests/auto/InvalidValue/2.31/good/Socket_PassCredentials_off.target
[Unit]
Description=Foo
[Socket]
PassCredentials=off

