# tests/auto/InvalidValue/2.31/good/Socket_LimitAS_0.target
[Unit]
Description=Foo
[Socket]
LimitAS=0

