# tests/auto/InvalidValue/2.31/good/Socket_PrivateNetwork_true.target
[Unit]
Description=Foo
[Socket]
PrivateNetwork=true

