# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EUCLEAN.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EUCLEAN

