# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EOVERFLOW.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EOVERFLOW

