# tests/auto/InvalidValue/2.31/good/Socket_SyslogLevelPrefix_no.target
[Unit]
Description=Foo
[Socket]
SyslogLevelPrefix=no

