# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGINFO.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGINFO

