# tests/auto/InvalidValue/2.31/good/Socket_TTYVTDisallocate_1.target
[Unit]
Description=Foo
[Socket]
TTYVTDisallocate=1

