# tests/auto/InvalidValue/2.31/good/Socket_PassSecurity_true.target
[Unit]
Description=Foo
[Socket]
PassSecurity=true

