# tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_27y.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=27y

