# tests/auto/InvalidValue/2.31/good/Socket_LimitNICE_9999999999999.target
[Unit]
Description=Foo
[Socket]
LimitNICE=9999999999999

