# tests/auto/InvalidValue/2.31/good/Socket_TTYReset_false.target
[Unit]
Description=Foo
[Socket]
TTYReset=false

