# tests/auto/InvalidValue/2.31/good/Socket_LimitCPU_9999999999999.target
[Unit]
Description=Foo
[Socket]
LimitCPU=9999999999999

