# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EMULTIHOP.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EMULTIHOP

