# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGSEGV.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGSEGV

