# tests/auto/InvalidValue/2.31/good/Socket_DisableControllers_blkio.target
[Unit]
Description=Foo
[Socket]
DisableControllers=blkio

