# tests/auto/InvalidValue/2.31/good/Socket_FreeBind_off.target
[Unit]
Description=Foo
[Socket]
FreeBind=off

