# tests/auto/InvalidValue/2.31/good/Socket_DisableControllers_cpu.target
[Unit]
Description=Foo
[Socket]
DisableControllers=cpu

