# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGALRM.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGALRM

