# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EADDRNOTAVAIL.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EADDRNOTAVAIL

