# tests/auto/InvalidValue/2.31/good/Socket_SecureBits_no-setuid-fixup.target
[Unit]
Description=Foo
[Socket]
SecureBits=no-setuid-fixup

