# tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_infinity.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=infinity

