# tests/auto/InvalidValue/2.31/good/Socket_RemoveIPC_on.target
[Unit]
Description=Foo
[Socket]
RemoveIPC=on

