# tests/auto/InvalidValue/2.31/good/Socket_TasksMax_infinity.target
[Unit]
Description=Foo
[Socket]
TasksMax=infinity

