# tests/auto/InvalidValue/2.31/good/Socket_SendSIGHUP_true.target
[Unit]
Description=Foo
[Socket]
SendSIGHUP=true

