# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGTERM.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGTERM

