# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ECONNREFUSED.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ECONNREFUSED

