# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGIO.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGIO

