# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGXCPU.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGXCPU

