# tests/auto/InvalidValue/2.31/good/Socket_SyslogFacility_local1.target
[Unit]
Description=Foo
[Socket]
SyslogFacility=local1

