# tests/auto/InvalidValue/2.31/good/Socket_TasksAccounting_yes.target
[Unit]
Description=Foo
[Socket]
TasksAccounting=yes

