# tests/auto/InvalidValue/2.31/good/Socket_Writable_true.target
[Unit]
Description=Foo
[Socket]
Writable=true

