# tests/auto/InvalidValue/2.31/good/Socket_CPUSchedulingResetOnFork_off.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingResetOnFork=off

