# tests/auto/InvalidValue/2.31/good/Socket_SyslogLevelPrefix_yes.target
[Unit]
Description=Foo
[Socket]
SyslogLevelPrefix=yes

