# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOTBLK.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOTBLK

