# tests/auto/InvalidValue/2.31/good/Socket_KillMode_none.target
[Unit]
Description=Foo
[Socket]
KillMode=none

