# tests/auto/InvalidValue/2.31/good/Socket_DisableControllers_pids.target
[Unit]
Description=Foo
[Socket]
DisableControllers=pids

