# tests/auto/InvalidValue/2.31/good/Socket_DirectoryMode_007.target
[Unit]
Description=Foo
[Socket]
DirectoryMode=007

