# tests/auto/InvalidValue/2.31/good/Socket_BindIPv6Only_ipv6-only.target
[Unit]
Description=Foo
[Socket]
BindIPv6Only=ipv6-only

