# tests/auto/InvalidValue/2.31/good/Socket_Personality_s390x.target
[Unit]
Description=Foo
[Socket]
Personality=s390x

