# tests/auto/InvalidValue/2.31/good/Socket_SELinuxContextFromNet_true.target
[Unit]
Description=Foo
[Socket]
SELinuxContextFromNet=true

