# tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitIntervalSec_15ms.target
[Unit]
Description=Foo
[Socket]
TriggerLimitIntervalSec=15ms

