# tests/auto/InvalidValue/2.31/good/Socket_KeepAliveIntervalSec_9999999999999.target
[Unit]
Description=Foo
[Socket]
KeepAliveIntervalSec=9999999999999

