# tests/auto/InvalidValue/2.31/good/Socket_TasksAccounting_no.target
[Unit]
Description=Foo
[Socket]
TasksAccounting=no

