# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_25.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=25

