# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGQUIT.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGQUIT

