# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EEXIST.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EEXIST

