# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EISDIR.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EISDIR

