# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EBADSLT.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EBADSLT

