# tests/auto/InvalidValue/2.31/good/Socket_LimitMSGQUEUE_9999999999999.target
[Unit]
Description=Foo
[Socket]
LimitMSGQUEUE=9999999999999

