# tests/auto/InvalidValue/2.31/good/Socket_LimitFSIZE_9999999999999.target
[Unit]
Description=Foo
[Socket]
LimitFSIZE=9999999999999

