# tests/auto/InvalidValue/2.31/good/Socket_Delegate_io.target
[Unit]
Description=Foo
[Socket]
Delegate=io

