# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EBADRQC.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EBADRQC

