# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_46.target
[Unit]
Description=Foo
[Socket]
KillSignal=46

