# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOKEY.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOKEY

