# tests/auto/InvalidValue/2.31/good/Socket_User_metoo.target
[Unit]
Description=Foo
[Socket]
User=metoo

