# tests/auto/InvalidValue/2.31/good/Socket_SendBuffer_0G.target
[Unit]
Description=Foo
[Socket]
SendBuffer=0G

