# tests/auto/InvalidValue/2.31/good/Socket_SELinuxContextFromNet_yes.target
[Unit]
Description=Foo
[Socket]
SELinuxContextFromNet=yes

