# tests/auto/InvalidValue/2.31/good/Socket_StandardOutput_abc.target
[Unit]
Description=Foo
[Socket]
StandardOutput=abc

