# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EXFULL.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EXFULL

