# tests/auto/InvalidValue/2.31/good/Socket_RootImage_--usr--sbin.target
[Unit]
Description=Foo
[Socket]
RootImage=/usr/sbin

