# tests/auto/InvalidValue/2.31/good/Socket_SyslogFacility_ftp.target
[Unit]
Description=Foo
[Socket]
SyslogFacility=ftp

