# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGHUP.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGHUP

