# tests/auto/InvalidValue/2.31/good/Socket_FreeBind_true.target
[Unit]
Description=Foo
[Socket]
FreeBind=true

