# tests/auto/InvalidValue/2.31/good/Socket_SocketGroup_metoo.target
[Unit]
Description=Foo
[Socket]
SocketGroup=metoo

