# tests/auto/InvalidValue/2.31/good/Socket_SendBuffer_9999999999999B.target
[Unit]
Description=Foo
[Socket]
SendBuffer=9999999999999B

