# tests/auto/InvalidValue/2.31/good/Socket_PrivateNetwork_0.target
[Unit]
Description=Foo
[Socket]
PrivateNetwork=0

