# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EFAULT.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EFAULT

