# tests/auto/InvalidValue/2.31/good/Socket_ReusePort_true.target
[Unit]
Description=Foo
[Socket]
ReusePort=true

