# tests/auto/InvalidValue/2.31/good/Socket_CPUSchedulingResetOnFork_0.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingResetOnFork=0

