# tests/auto/InvalidValue/2.31/good/Socket_MountFlags_shared.target
[Unit]
Description=Foo
[Socket]
MountFlags=shared

