# tests/auto/InvalidValue/2.31/good/Socket_SyslogFacility_daemon.target
[Unit]
Description=Foo
[Socket]
SyslogFacility=daemon

