# tests/auto/InvalidValue/2.31/good/Socket_MountAPIVFS_on.target
[Unit]
Description=Foo
[Socket]
MountAPIVFS=on

