# tests/auto/InvalidValue/2.31/good/Socket_WorkingDirectory_--bin.target
[Unit]
Description=Foo
[Socket]
WorkingDirectory=/bin

