# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOSTR.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOSTR

