# tests/auto/InvalidValue/2.31/good/Socket_MemoryMax_0K.target
[Unit]
Description=Foo
[Socket]
MemoryMax=0K

