# tests/auto/InvalidValue/2.31/good/Socket_SyslogLevel_notice.target
[Unit]
Description=Foo
[Socket]
SyslogLevel=notice

