# tests/auto/InvalidValue/2.31/good/Socket_SendSIGKILL_true.target
[Unit]
Description=Foo
[Socket]
SendSIGKILL=true

