# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_24.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=24

