# tests/auto/InvalidValue/2.31/good/Socket_ReusePort_false.target
[Unit]
Description=Foo
[Socket]
ReusePort=false

