# tests/auto/InvalidValue/2.31/good/Socket_FreeBind_yes.target
[Unit]
Description=Foo
[Socket]
FreeBind=yes

