# tests/auto/InvalidValue/2.31/good/Socket_ListenSpecial_--dev--.target
[Unit]
Description=Foo
[Socket]
ListenSpecial=/dev/

