# tests/auto/InvalidValue/2.31/good/Socket_Broadcast_false.target
[Unit]
Description=Foo
[Socket]
Broadcast=false

