# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_55.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=55

