# tests/auto/InvalidValue/2.31/good/Socket_SocketMode_777.target
[Unit]
Description=Foo
[Socket]
SocketMode=777

