# tests/auto/InvalidValue/2.31/good/Socket_PassEnvironment_abc.target
[Unit]
Description=Foo
[Socket]
PassEnvironment=abc

