# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGHUP.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGHUP

