# tests/auto/InvalidValue/2.31/good/Socket_TasksAccounting_false.target
[Unit]
Description=Foo
[Socket]
TasksAccounting=false

