# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ESRCH.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ESRCH

