# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOTTY.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOTTY

