# tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitIntervalSec_17second.target
[Unit]
Description=Foo
[Socket]
TriggerLimitIntervalSec=17second

