# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_46.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=46

