# tests/auto/InvalidValue/2.31/good/Socket_IOAccounting_yes.target
[Unit]
Description=Foo
[Socket]
IOAccounting=yes

