# tests/auto/InvalidValue/2.31/good/Socket_Writable_on.target
[Unit]
Description=Foo
[Socket]
Writable=on

