# tests/auto/InvalidValue/2.31/good/Socket_BlockStartupIOWeight_1.target
[Unit]
Description=Foo
[Socket]
BlockStartupIOWeight=1

