# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGALRM.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGALRM

