# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGQUIT.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGQUIT

