# tests/auto/InvalidValue/2.31/good/Socket_PrivateTmp_1.target
[Unit]
Description=Foo
[Socket]
PrivateTmp=1

