# tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_23week.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=23week

