# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EMFILE.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EMFILE

