# tests/auto/InvalidValue/2.31/good/Socket_MessageQueueMaxMessages_0.target
[Unit]
Description=Foo
[Socket]
MessageQueueMaxMessages=0

