# tests/auto/InvalidValue/2.31/good/Socket_IPTOS_9999999999999.target
[Unit]
Description=Foo
[Socket]
IPTOS=9999999999999

