# tests/auto/InvalidValue/2.31/good/Socket_UtmpMode_user.target
[Unit]
Description=Foo
[Socket]
UtmpMode=user

