# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGXCPU.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGXCPU

