# tests/auto/InvalidValue/2.31/good/Socket_TTYReset_true.target
[Unit]
Description=Foo
[Socket]
TTYReset=true

