# tests/auto/InvalidValue/2.31/good/Socket_Broadcast_off.target
[Unit]
Description=Foo
[Socket]
Broadcast=off

