# tests/auto/InvalidValue/2.31/good/Socket_SyslogLevel_err.target
[Unit]
Description=Foo
[Socket]
SyslogLevel=err

