# tests/auto/InvalidValue/2.31/good/Socket_LimitNICE_infinity.target
[Unit]
Description=Foo
[Socket]
LimitNICE=infinity

