# tests/auto/InvalidValue/2.31/good/Socket_CPUShares_262144.target
[Unit]
Description=Foo
[Socket]
CPUShares=262144

