# tests/auto/InvalidValue/2.31/good/Socket_SyslogFacility_auth.target
[Unit]
Description=Foo
[Socket]
SyslogFacility=auth

