# tests/auto/InvalidValue/2.31/good/Socket_StartupCPUWeight_1.target
[Unit]
Description=Foo
[Socket]
StartupCPUWeight=1

