# tests/auto/InvalidValue/2.31/good/Socket_RemoveIPC_1.target
[Unit]
Description=Foo
[Socket]
RemoveIPC=1

