# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGSEGV.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGSEGV

