# tests/auto/InvalidValue/2.31/good/Socket_IPTOS_low-cost.target
[Unit]
Description=Foo
[Socket]
IPTOS=low-cost

