# tests/auto/InvalidValue/2.31/good/Socket_Writable_yes.target
[Unit]
Description=Foo
[Socket]
Writable=yes

