# tests/auto/InvalidValue/2.31/good/Socket_SELinuxContextFromNet_no.target
[Unit]
Description=Foo
[Socket]
SELinuxContextFromNet=no

