# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGLOST.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGLOST

