# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGTTOU.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGTTOU

