# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EINTR.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EINTR

