# tests/auto/InvalidValue/2.31/good/Socket_IPTTL_0.target
[Unit]
Description=Foo
[Socket]
IPTTL=0

