# tests/auto/InvalidValue/2.31/good/Socket_SyslogFacility_cron.target
[Unit]
Description=Foo
[Socket]
SyslogFacility=cron

