# tests/auto/InvalidValue/2.31/good/Socket_UMask_9999999999999.target
[Unit]
Description=Foo
[Socket]
UMask=9999999999999

