# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_40.target
[Unit]
Description=Foo
[Socket]
KillSignal=40

