# tests/auto/InvalidValue/2.31/good/Socket_SendSIGKILL_no.target
[Unit]
Description=Foo
[Socket]
SendSIGKILL=no

