# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ELNRANGE.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ELNRANGE

