# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGPIPE.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGPIPE

