# tests/auto/InvalidValue/2.31/good/Socket_TTYVHangup_0.target
[Unit]
Description=Foo
[Socket]
TTYVHangup=0

