# tests/auto/InvalidValue/2.31/good/Socket_PrivateDevices_0.target
[Unit]
Description=Foo
[Socket]
PrivateDevices=0

