# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_34.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=34

