# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EXDEV.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EXDEV

