# tests/auto/InvalidValue/2.31/good/Socket_StandardError_abc.target
[Unit]
Description=Foo
[Socket]
StandardError=abc

