# tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_22weeks.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=22weeks

