# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGPIPE.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGPIPE

