# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENAMETOOLONG.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENAMETOOLONG

