# tests/auto/InvalidValue/2.31/good/Socket_TasksAccounting_true.target
[Unit]
Description=Foo
[Socket]
TasksAccounting=true

