# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_11.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=11

