# tests/auto/InvalidValue/2.31/good/Socket_KillMode_process.target
[Unit]
Description=Foo
[Socket]
KillMode=process

