# tests/auto/InvalidValue/2.31/good/Socket_SyslogLevelPrefix_0.target
[Unit]
Description=Foo
[Socket]
SyslogLevelPrefix=0

