# tests/auto/InvalidValue/2.31/good/Socket_FreeBind_false.target
[Unit]
Description=Foo
[Socket]
FreeBind=false

