# tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitIntervalSec_4hour.target
[Unit]
Description=Foo
[Socket]
TriggerLimitIntervalSec=4hour

