# tests/auto/InvalidValue/2.31/good/Socket_MountAPIVFS_true.target
[Unit]
Description=Foo
[Socket]
MountAPIVFS=true

