# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_1.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=1

