# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ETXTBSY.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ETXTBSY

