# tests/auto/InvalidValue/2.31/good/Socket_StartupCPUShares_262144.target
[Unit]
Description=Foo
[Socket]
StartupCPUShares=262144

