# tests/auto/InvalidValue/2.31/good/Socket_LimitDATA_9999999999999.target
[Unit]
Description=Foo
[Socket]
LimitDATA=9999999999999

