# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EOPNOTSUPP.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EOPNOTSUPP

