# tests/auto/InvalidValue/2.31/good/Socket_TasksAccounting_on.target
[Unit]
Description=Foo
[Socket]
TasksAccounting=on

