# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGTRAP.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGTRAP

