# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EL2NSYNC.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EL2NSYNC

