# tests/auto/InvalidValue/2.31/good/Socket_IOReadIOPSMax_abc.target
[Unit]
Description=Foo
[Socket]
IOReadIOPSMax=abc

