# tests/auto/InvalidValue/2.31/good/Socket_LimitNOFILE_infinity.target
[Unit]
Description=Foo
[Socket]
LimitNOFILE=infinity

