# tests/auto/InvalidValue/2.31/good/Socket_SyslogLevel_debug.target
[Unit]
Description=Foo
[Socket]
SyslogLevel=debug

