# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EBADMSG.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EBADMSG

