# tests/auto/InvalidValue/2.31/good/Socket_ProtectSystem_strict.target
[Unit]
Description=Foo
[Socket]
ProtectSystem=strict

