# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EKEYREJECTED.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EKEYREJECTED

