# tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitIntervalSec_3hours.target
[Unit]
Description=Foo
[Socket]
TriggerLimitIntervalSec=3hours

