# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGIOT.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGIOT

