# tests/auto/InvalidValue/2.31/good/Socket_ProtectHome_off.target
[Unit]
Description=Foo
[Socket]
ProtectHome=off

