# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENETDOWN.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENETDOWN

