# tests/auto/InvalidValue/2.31/good/Socket_Backlog_9999999999999.target
[Unit]
Description=Foo
[Socket]
Backlog=9999999999999

