# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ELIBBAD.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ELIBBAD

