# tests/auto/InvalidValue/2.31/good/Socket_PassCredentials_no.target
[Unit]
Description=Foo
[Socket]
PassCredentials=no

