# tests/auto/InvalidValue/2.31/good/Socket_MemoryDenyWriteExecute_1.target
[Unit]
Description=Foo
[Socket]
MemoryDenyWriteExecute=1

