# tests/auto/InvalidValue/2.31/good/Socket_Writable_no.target
[Unit]
Description=Foo
[Socket]
Writable=no

