# tests/auto/InvalidValue/2.31/good/Socket_CPUSchedulingPolicy_batch.target
[Unit]
Description=Foo
[Socket]
CPUSchedulingPolicy=batch

