# tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitIntervalSec_11m.target
[Unit]
Description=Foo
[Socket]
TriggerLimitIntervalSec=11m

