# tests/auto/InvalidValue/2.31/good/Socket_SmackLabelIPOut_abc.target
[Unit]
Description=Foo
[Socket]
SmackLabelIPOut=abc

