# tests/auto/InvalidValue/2.31/good/Socket_UtmpMode_init.target
[Unit]
Description=Foo
[Socket]
UtmpMode=init

