# tests/auto/InvalidValue/2.31/good/Socket_TTYVTDisallocate_false.target
[Unit]
Description=Foo
[Socket]
TTYVTDisallocate=false

