# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGIOT.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGIOT

