# tests/auto/InvalidValue/2.31/good/Socket_SendSIGKILL_0.target
[Unit]
Description=Foo
[Socket]
SendSIGKILL=0

