# tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_13months.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=13months

