# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGABRT.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGABRT

