# tests/auto/InvalidValue/2.31/good/Socket_MountAPIVFS_1.target
[Unit]
Description=Foo
[Socket]
MountAPIVFS=1

