# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EINPROGRESS.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EINPROGRESS

