# tests/auto/InvalidValue/2.31/good/Socket_CapabilityBoundingSet_abc.target
[Unit]
Description=Foo
[Socket]
CapabilityBoundingSet=abc

