# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ECANCELED.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ECANCELED

