# tests/auto/InvalidValue/2.31/good/Socket_PrivateNetwork_false.target
[Unit]
Description=Foo
[Socket]
PrivateNetwork=false

