# tests/auto/InvalidValue/2.31/good/Socket_RemoveIPC_true.target
[Unit]
Description=Foo
[Socket]
RemoveIPC=true

