# tests/auto/InvalidValue/2.31/good/Socket_DevicePolicy_auto.target
[Unit]
Description=Foo
[Socket]
DevicePolicy=auto

