# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGWINCH.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGWINCH

