# tests/auto/InvalidValue/2.31/good/Socket_LimitFSIZE_infinity.target
[Unit]
Description=Foo
[Socket]
LimitFSIZE=infinity

