# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOLCK.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOLCK

