# tests/auto/InvalidValue/2.31/good/Socket_StartupIOWeight_1.target
[Unit]
Description=Foo
[Socket]
StartupIOWeight=1

