# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_3.target
[Unit]
Description=Foo
[Socket]
KillSignal=3

