# tests/auto/InvalidValue/2.31/good/Socket_MemoryLow_9999999999999K.target
[Unit]
Description=Foo
[Socket]
MemoryLow=9999999999999K

