# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGTTOU.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGTTOU

