# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_42.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=42

