# tests/auto/InvalidValue/2.31/good/Socket_MountFlags_slave.target
[Unit]
Description=Foo
[Socket]
MountFlags=slave

