# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EREMCHG.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EREMCHG

