# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGINT.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGINT

