# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ELIBACC.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ELIBACC

