# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOPKG.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOPKG

