# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EPIPE.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EPIPE

