# tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitIntervalSec_12month.target
[Unit]
Description=Foo
[Socket]
TriggerLimitIntervalSec=12month

