# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGIOT.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGIOT

