# tests/auto/InvalidValue/2.31/good/Socket_IOSchedulingPriority_7.target
[Unit]
Description=Foo
[Socket]
IOSchedulingPriority=7

