# tests/auto/InvalidValue/2.31/good/Socket_PassCredentials_1.target
[Unit]
Description=Foo
[Socket]
PassCredentials=1

