# tests/auto/InvalidValue/2.31/good/Socket_TimeoutSec_18sec.target
[Unit]
Description=Foo
[Socket]
TimeoutSec=18sec

