# tests/auto/InvalidValue/2.31/good/Socket_RestrictRealtime_1.target
[Unit]
Description=Foo
[Socket]
RestrictRealtime=1

