# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_65.target
[Unit]
Description=Foo
[Socket]
KillSignal=65

