# tests/auto/InvalidValue/2.31/good/Socket_OOMScoreAdjust_-1000.target
[Unit]
Description=Foo
[Socket]
OOMScoreAdjust=-1000

