# tests/auto/InvalidValue/2.31/good/Socket_FreeBind_on.target
[Unit]
Description=Foo
[Socket]
FreeBind=on

