# tests/auto/InvalidValue/2.31/good/Socket_RemoveOnStop_false.target
[Unit]
Description=Foo
[Socket]
RemoveOnStop=false

