# tests/auto/InvalidValue/2.31/good/Socket_MemoryMax_9999999999999T.target
[Unit]
Description=Foo
[Socket]
MemoryMax=9999999999999T

