# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGCLD.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGCLD

