# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGTERM.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGTERM

