# tests/auto/InvalidValue/2.31/good/Socket_ExecStopPost_--bin--sh.target
[Unit]
Description=Foo
[Socket]
ExecStopPost=/bin/sh

