# tests/auto/InvalidValue/2.31/good/Socket_CPUAffinity_numa.target
[Unit]
Description=Foo
[Socket]
CPUAffinity=numa

