# tests/auto/InvalidValue/2.31/good/Socket_IODeviceWeight_abc.target
[Unit]
Description=Foo
[Socket]
IODeviceWeight=abc

