# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGPWR.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGPWR

