# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGEMT.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGEMT

