# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EIO.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EIO

