# tests/auto/InvalidValue/2.31/good/Socket_TTYVTDisallocate_off.target
[Unit]
Description=Foo
[Socket]
TTYVTDisallocate=off

