# tests/auto/InvalidValue/2.31/good/Socket_SocketUser_metoo.target
[Unit]
Description=Foo
[Socket]
SocketUser=metoo

