# tests/auto/InvalidValue/2.31/good/Socket_SystemCallArchitectures_abc.target
[Unit]
Description=Foo
[Socket]
SystemCallArchitectures=abc

