# tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_20usec.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=20usec

