# tests/auto/InvalidValue/2.31/good/Socket_StandardOutput_null.target
[Unit]
Description=Foo
[Socket]
StandardOutput=null

