# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOSPC.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOSPC

