# tests/auto/InvalidValue/2.31/good/Socket_SocketProtocol_udplite.target
[Unit]
Description=Foo
[Socket]
SocketProtocol=udplite

