# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ECOMM.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ECOMM

