# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGKILL.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGKILL

