# tests/auto/InvalidValue/2.31/good/Socket_Delegate_devices.target
[Unit]
Description=Foo
[Socket]
Delegate=devices

