# tests/auto/InvalidValue/2.31/good/Socket_RemoveOnStop_yes.target
[Unit]
Description=Foo
[Socket]
RemoveOnStop=yes

