# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_51.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=51

