# tests/auto/InvalidValue/2.31/good/Socket_Priority_9999999999999.target
[Unit]
Description=Foo
[Socket]
Priority=9999999999999

