# tests/auto/InvalidValue/2.31/good/Socket_ListenMessageQueue_abc.target
[Unit]
Description=Foo
[Socket]
ListenMessageQueue=abc

