# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGTSTP.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGTSTP

