# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EINVAL.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EINVAL

