# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EPFNOSUPPORT.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EPFNOSUPPORT

