# tests/auto/InvalidValue/2.31/good/Socket_SyslogLevelPrefix_false.target
[Unit]
Description=Foo
[Socket]
SyslogLevelPrefix=false

