# tests/auto/InvalidValue/2.31/good/Socket_Transparent_false.target
[Unit]
Description=Foo
[Socket]
Transparent=false

