# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_59.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=59

