# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGLOST.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGLOST

