# tests/auto/InvalidValue/2.31/good/Socket_LimitCPU_infinity.target
[Unit]
Description=Foo
[Socket]
LimitCPU=infinity

