# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EILSEQ.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EILSEQ

