# tests/auto/InvalidValue/2.31/good/Socket_TasksMax_1.target
[Unit]
Description=Foo
[Socket]
TasksMax=1

