# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EKEYEXPIRED.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EKEYEXPIRED

