# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGPOLL.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGPOLL

