# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_41.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=41

