# tests/auto/InvalidValue/2.31/good/Socket_TasksAccounting_1.target
[Unit]
Description=Foo
[Socket]
TasksAccounting=1

