# tests/auto/InvalidValue/2.31/good/Socket_IOSchedulingClass_none.target
[Unit]
Description=Foo
[Socket]
IOSchedulingClass=none

