# tests/auto/InvalidValue/2.31/good/Socket_ListenDatagram_abc.target
[Unit]
Description=Foo
[Socket]
ListenDatagram=abc

