# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENXIO.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENXIO

