# tests/auto/InvalidValue/2.31/good/Socket_SyslogLevelPrefix_on.target
[Unit]
Description=Foo
[Socket]
SyslogLevelPrefix=on

