# tests/auto/InvalidValue/2.31/good/Socket_TTYReset_0.target
[Unit]
Description=Foo
[Socket]
TTYReset=0

