# tests/auto/InvalidValue/2.31/good/Socket_PrivateTmp_no.target
[Unit]
Description=Foo
[Socket]
PrivateTmp=no

