# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGSYS.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGSYS

