# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EL3HLT.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EL3HLT

