# tests/auto/InvalidValue/2.31/good/Socket_IOAccounting_no.target
[Unit]
Description=Foo
[Socket]
IOAccounting=no

