# tests/auto/InvalidValue/2.31/good/Socket_ListenStream_abc.target
[Unit]
Description=Foo
[Socket]
ListenStream=abc

