# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGIO.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGIO

