# tests/auto/InvalidValue/2.31/good/Socket_KeepAlive_false.target
[Unit]
Description=Foo
[Socket]
KeepAlive=false

