# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGTTIN.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGTTIN

