# tests/auto/InvalidValue/2.31/good/Socket_SocketGroup_iamagroup.target
[Unit]
Description=Foo
[Socket]
SocketGroup=iamagroup

