# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_SIGILL.target
[Unit]
Description=Foo
[Socket]
KillSignal=SIGILL

