# tests/auto/InvalidValue/2.31/good/Socket_FreeBind_no.target
[Unit]
Description=Foo
[Socket]
FreeBind=no

