# tests/auto/InvalidValue/2.31/good/Socket_LimitAS_9999999999999.target
[Unit]
Description=Foo
[Socket]
LimitAS=9999999999999

