# tests/auto/InvalidValue/2.31/good/Socket_Personality_ppc.target
[Unit]
Description=Foo
[Socket]
Personality=ppc

