# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOTUNIQ.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOTUNIQ

