# tests/auto/InvalidValue/2.31/good/Socket_RestrictRealtime_yes.target
[Unit]
Description=Foo
[Socket]
RestrictRealtime=yes

