# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOANO.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOANO

