# tests/auto/InvalidValue/2.31/good/Socket_SELinuxContextFromNet_false.target
[Unit]
Description=Foo
[Socket]
SELinuxContextFromNet=false

