# tests/auto/InvalidValue/2.31/good/Socket_RemoveOnStop_on.target
[Unit]
Description=Foo
[Socket]
RemoveOnStop=on

