# tests/auto/InvalidValue/2.31/good/Socket_Writable_false.target
[Unit]
Description=Foo
[Socket]
Writable=false

