# tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_25years.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=25years

