# tests/auto/InvalidValue/2.31/good/Socket_LimitSTACK_0.target
[Unit]
Description=Foo
[Socket]
LimitSTACK=0

