# tests/auto/InvalidValue/2.31/good/Socket_LimitMSGQUEUE_0.target
[Unit]
Description=Foo
[Socket]
LimitMSGQUEUE=0

