# tests/auto/InvalidValue/2.31/good/Socket_LimitRTTIME_0.target
[Unit]
Description=Foo
[Socket]
LimitRTTIME=0

