# tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitIntervalSec_23week.target
[Unit]
Description=Foo
[Socket]
TriggerLimitIntervalSec=23week

