# tests/auto/InvalidValue/2.31/good/Socket_IOSchedulingClass_realtime.target
[Unit]
Description=Foo
[Socket]
IOSchedulingClass=realtime

