# tests/auto/InvalidValue/2.31/good/Socket_ReusePort_yes.target
[Unit]
Description=Foo
[Socket]
ReusePort=yes

