# tests/auto/InvalidValue/2.31/good/Socket_PassSecurity_off.target
[Unit]
Description=Foo
[Socket]
PassSecurity=off

