# tests/auto/InvalidValue/2.31/good/Socket_LimitRTTIME_9999999999999.target
[Unit]
Description=Foo
[Socket]
LimitRTTIME=9999999999999

