# tests/auto/InvalidValue/2.31/good/Socket_StandardError_journal%2Bconsole.target
[Unit]
Description=Foo
[Socket]
StandardError=journal+console

