# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOLINK.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOLINK

