# tests/auto/InvalidValue/2.31/good/Socket_RestrictSUIDSGID_0.target
[Unit]
Description=Foo
[Socket]
RestrictSUIDSGID=0

