# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EPERM.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EPERM

