# tests/auto/InvalidValue/2.31/good/Socket_ReadWriteDirectories_abc.target
[Unit]
Description=Foo
[Socket]
ReadWriteDirectories=abc

