# tests/auto/InvalidValue/2.31/good/Socket_Group_iamagroup.target
[Unit]
Description=Foo
[Socket]
Group=iamagroup

