# tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_12month.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=12month

