# tests/auto/InvalidValue/2.31/good/Socket_AppArmorProfile_abc.target
[Unit]
Description=Foo
[Socket]
AppArmorProfile=abc

