# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGSYS.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGSYS

