# tests/auto/InvalidValue/2.31/good/Socket_MemoryDenyWriteExecute_false.target
[Unit]
Description=Foo
[Socket]
MemoryDenyWriteExecute=false

