# tests/auto/InvalidValue/2.31/good/Socket_PrivateDevices_on.target
[Unit]
Description=Foo
[Socket]
PrivateDevices=on

