# tests/auto/InvalidValue/2.31/good/Socket_KillSignal_17.target
[Unit]
Description=Foo
[Socket]
KillSignal=17

