# tests/auto/InvalidValue/2.31/good/Socket_Transparent_yes.target
[Unit]
Description=Foo
[Socket]
Transparent=yes

