# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_33.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=33

