# tests/auto/InvalidValue/2.31/good/Socket_TriggerLimitBurst_21us.target
[Unit]
Description=Foo
[Socket]
TriggerLimitBurst=21us

