# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EBADE.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EBADE

