# tests/auto/InvalidValue/2.31/good/Socket_IgnoreSIGPIPE_no.target
[Unit]
Description=Foo
[Socket]
IgnoreSIGPIPE=no

