# tests/auto/InvalidValue/2.31/good/Socket_ProtectHome_false.target
[Unit]
Description=Foo
[Socket]
ProtectHome=false

