# tests/auto/InvalidValue/2.31/good/Socket_ExecStartPost_--bin--sh.target
[Unit]
Description=Foo
[Socket]
ExecStartPost=/bin/sh

