# tests/auto/InvalidValue/2.31/good/Socket_ProtectHome_tmpfs.target
[Unit]
Description=Foo
[Socket]
ProtectHome=tmpfs

