# tests/auto/InvalidValue/2.31/good/Socket_ReusePort_1.target
[Unit]
Description=Foo
[Socket]
ReusePort=1

