# tests/auto/InvalidValue/2.31/good/Socket_RestrictSUIDSGID_true.target
[Unit]
Description=Foo
[Socket]
RestrictSUIDSGID=true

