# tests/auto/InvalidValue/2.31/good/Socket_ListenSequentialPacket_abc.target
[Unit]
Description=Foo
[Socket]
ListenSequentialPacket=abc

