# tests/auto/InvalidValue/2.31/good/Socket_WatchdogSignal_SIGBUS.target
[Unit]
Description=Foo
[Socket]
WatchdogSignal=SIGBUS

