# tests/auto/InvalidValue/2.31/good/Socket_MountAPIVFS_0.target
[Unit]
Description=Foo
[Socket]
MountAPIVFS=0

