# tests/auto/InvalidValue/2.31/good/Socket_ReadWritePaths_abc.target
[Unit]
Description=Foo
[Socket]
ReadWritePaths=abc

