# tests/auto/InvalidValue/2.31/good/Socket_PassCredentials_yes.target
[Unit]
Description=Foo
[Socket]
PassCredentials=yes

