# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EMEDIUMTYPE.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EMEDIUMTYPE

