# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EBUSY.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EBUSY

