# tests/auto/InvalidValue/2.31/good/Socket_PrivateNetwork_on.target
[Unit]
Description=Foo
[Socket]
PrivateNetwork=on

