# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EADDRINUSE.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EADDRINUSE

