# tests/auto/InvalidValue/2.31/good/Socket_IOAccounting_on.target
[Unit]
Description=Foo
[Socket]
IOAccounting=on

