# tests/auto/InvalidValue/2.31/good/Socket_IgnoreSIGPIPE_on.target
[Unit]
Description=Foo
[Socket]
IgnoreSIGPIPE=on

