# tests/auto/InvalidValue/2.31/good/Socket_PassCredentials_true.target
[Unit]
Description=Foo
[Socket]
PassCredentials=true

