# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EAFNOSUPPORT.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EAFNOSUPPORT

