# tests/auto/InvalidValue/2.31/good/Socket_MountFlags_private.target
[Unit]
Description=Foo
[Socket]
MountFlags=private

