# tests/auto/InvalidValue/2.31/good/Socket_RootDirectory_--usr--sbin.target
[Unit]
Description=Foo
[Socket]
RootDirectory=/usr/sbin

