# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EREMOTE.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EREMOTE

