# tests/auto/InvalidValue/2.31/good/Socket_TimerSlackNSec_9999999999999.target
[Unit]
Description=Foo
[Socket]
TimerSlackNSec=9999999999999

