# tests/auto/InvalidValue/2.31/good/Socket_FinalKillSignal_SIGILL.target
[Unit]
Description=Foo
[Socket]
FinalKillSignal=SIGILL

