# tests/auto/InvalidValue/2.31/good/Socket_Broadcast_0.target
[Unit]
Description=Foo
[Socket]
Broadcast=0

