# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ENOTSUP.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ENOTSUP

