# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_ESOCKTNOSUPPORT.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=ESOCKTNOSUPPORT

