# tests/auto/InvalidValue/2.31/good/Socket_SyslogFacility_kern.target
[Unit]
Description=Foo
[Socket]
SyslogFacility=kern

