# tests/auto/InvalidValue/2.31/good/Socket_SyslogLevel_alert.target
[Unit]
Description=Foo
[Socket]
SyslogLevel=alert

