# tests/auto/InvalidValue/2.31/good/Socket_TasksAccounting_0.target
[Unit]
Description=Foo
[Socket]
TasksAccounting=0

