# tests/auto/InvalidValue/2.31/good/Socket_SystemCallErrorNumber_EMLINK.target
[Unit]
Description=Foo
[Socket]
SystemCallErrorNumber=EMLINK

