# tests/auto/InvalidValue/2.31/good/Socket_TimeoutSec_14msec.target
[Unit]
Description=Foo
[Socket]
TimeoutSec=14msec

