# tests/auto/InvalidValue/2.31/good/Socket_Delegate_cpuacct.target
[Unit]
Description=Foo
[Socket]
Delegate=cpuacct

