# tests/auto/InvalidValue/2.31/good/Swap_RestrictRealtime_on.target
[Unit]
Description=Foo
[Swap]
RestrictRealtime=on

