# tests/auto/InvalidValue/2.31/good/Swap_LimitRSS_0.target
[Unit]
Description=Foo
[Swap]
LimitRSS=0

