# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EDEADLOCK.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EDEADLOCK

