# tests/auto/InvalidValue/2.31/good/Swap_PrivateDevices_on.target
[Unit]
Description=Foo
[Swap]
PrivateDevices=on

