# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGIOT.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGIOT

