# tests/auto/InvalidValue/2.31/good/Swap_MemoryAccounting_on.target
[Unit]
Description=Foo
[Swap]
MemoryAccounting=on

