# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EPROTO.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EPROTO

