# tests/auto/InvalidValue/2.31/good/Swap_MemoryDenyWriteExecute_false.target
[Unit]
Description=Foo
[Swap]
MemoryDenyWriteExecute=false

