# tests/auto/InvalidValue/2.31/good/Swap_IOAccounting_no.target
[Unit]
Description=Foo
[Swap]
IOAccounting=no

