# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ELIBMAX.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ELIBMAX

