# tests/auto/InvalidValue/2.31/good/Swap_SendSIGHUP_true.target
[Unit]
Description=Foo
[Swap]
SendSIGHUP=true

