# tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGLOST.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGLOST

