# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGXFSZ.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGXFSZ

