# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ELIBBAD.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ELIBBAD

