# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENODATA.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENODATA

