# tests/auto/InvalidValue/2.31/good/Swap_CapabilityBoundingSet_abc.target
[Unit]
Description=Foo
[Swap]
CapabilityBoundingSet=abc

