# tests/auto/InvalidValue/2.31/good/Swap_BlockIOAccounting_on.target
[Unit]
Description=Foo
[Swap]
BlockIOAccounting=on

