# tests/auto/InvalidValue/2.31/good/Swap_OOMScoreAdjust_-1000.target
[Unit]
Description=Foo
[Swap]
OOMScoreAdjust=-1000

