# tests/auto/InvalidValue/2.31/good/Swap_IOAccounting_false.target
[Unit]
Description=Foo
[Swap]
IOAccounting=false

