# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGSTOP.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGSTOP

