# tests/auto/InvalidValue/2.31/good/Swap_SyslogFacility_local3.target
[Unit]
Description=Foo
[Swap]
SyslogFacility=local3

