# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ELIBACC.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ELIBACC

