# tests/auto/InvalidValue/2.31/good/Swap_CPUAffinity_abc.target
[Unit]
Description=Foo
[Swap]
CPUAffinity=abc

