# tests/auto/InvalidValue/2.31/good/Swap_CPUSchedulingResetOnFork_on.target
[Unit]
Description=Foo
[Swap]
CPUSchedulingResetOnFork=on

