# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ETXTBSY.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ETXTBSY

