# tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGPOLL.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGPOLL

