# tests/auto/InvalidValue/2.31/good/Swap_CPUSchedulingResetOnFork_true.target
[Unit]
Description=Foo
[Swap]
CPUSchedulingResetOnFork=true

