# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ESHUTDOWN.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ESHUTDOWN

