# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_49.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=49

