# tests/auto/InvalidValue/2.31/good/Swap_UtmpMode_user.target
[Unit]
Description=Foo
[Swap]
UtmpMode=user

