# tests/auto/InvalidValue/2.31/good/Swap_PrivateTmp_false.target
[Unit]
Description=Foo
[Swap]
PrivateTmp=false

