# tests/auto/InvalidValue/2.31/good/Swap_RemoveIPC_0.target
[Unit]
Description=Foo
[Swap]
RemoveIPC=0

