# tests/auto/InvalidValue/2.31/good/Swap_KillMode_process.target
[Unit]
Description=Foo
[Swap]
KillMode=process

