# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOLCK.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOLCK

