# tests/auto/InvalidValue/2.31/good/Swap_IODeviceWeight_abc.target
[Unit]
Description=Foo
[Swap]
IODeviceWeight=abc

