# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EUNATCH.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EUNATCH

