# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ECONNREFUSED.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ECONNREFUSED

