# tests/auto/InvalidValue/2.31/good/Swap_TTYPath_--dev--foo.target
[Unit]
Description=Foo
[Swap]
TTYPath=/dev/foo

