# tests/auto/InvalidValue/2.31/good/Swap_SyslogFacility_uucp.target
[Unit]
Description=Foo
[Swap]
SyslogFacility=uucp

