# tests/auto/InvalidValue/2.31/good/Swap_RemoveIPC_yes.target
[Unit]
Description=Foo
[Swap]
RemoveIPC=yes

