# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGTTOU.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGTTOU

