# tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_13months.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=13months

