# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EIO.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EIO

