# tests/auto/InvalidValue/2.31/good/Swap_StandardOutput_socket.target
[Unit]
Description=Foo
[Swap]
StandardOutput=socket

