# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EBADSLT.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EBADSLT

