# tests/auto/InvalidValue/2.31/good/Swap_TTYReset_no.target
[Unit]
Description=Foo
[Swap]
TTYReset=no

