# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ETIME.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ETIME

