# tests/auto/InvalidValue/2.31/good/Swap_SendSIGKILL_1.target
[Unit]
Description=Foo
[Swap]
SendSIGKILL=1

