# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGINT.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGINT

