# tests/auto/InvalidValue/2.31/good/Swap_KillMode_mixed.target
[Unit]
Description=Foo
[Swap]
KillMode=mixed

