# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ECANCELED.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ECANCELED

