# tests/auto/InvalidValue/2.31/good/Swap_MountFlags_shared.target
[Unit]
Description=Foo
[Swap]
MountFlags=shared

