# tests/auto/InvalidValue/2.31/good/Swap_PrivateDevices_false.target
[Unit]
Description=Foo
[Swap]
PrivateDevices=false

