# tests/auto/InvalidValue/2.31/good/Swap_IOReadIOPSMax_abc.target
[Unit]
Description=Foo
[Swap]
IOReadIOPSMax=abc

