# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EBADE.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EBADE

