# tests/auto/InvalidValue/2.31/good/Swap_BlockIOAccounting_false.target
[Unit]
Description=Foo
[Swap]
BlockIOAccounting=false

