# tests/auto/InvalidValue/2.31/good/Swap_SyslogLevelPrefix_0.target
[Unit]
Description=Foo
[Swap]
SyslogLevelPrefix=0

