# tests/auto/InvalidValue/2.31/good/Swap_PrivateNetwork_on.target
[Unit]
Description=Foo
[Swap]
PrivateNetwork=on

