# tests/auto/InvalidValue/2.31/good/Swap_RestrictRealtime_1.target
[Unit]
Description=Foo
[Swap]
RestrictRealtime=1

