# tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGFPE.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGFPE

