# tests/auto/InvalidValue/2.31/good/Swap_PrivateTmp_0.target
[Unit]
Description=Foo
[Swap]
PrivateTmp=0

