# tests/auto/InvalidValue/2.31/good/Swap_PrivateTmp_yes.target
[Unit]
Description=Foo
[Swap]
PrivateTmp=yes

