# tests/auto/InvalidValue/2.31/good/Swap_Personality_ppc.target
[Unit]
Description=Foo
[Swap]
Personality=ppc

