# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_28.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=28

