# tests/auto/InvalidValue/2.31/good/Swap_SyslogLevelPrefix_1.target
[Unit]
Description=Foo
[Swap]
SyslogLevelPrefix=1

