# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EHOSTDOWN.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EHOSTDOWN

