# tests/auto/InvalidValue/2.31/good/Swap_SyslogLevel_notice.target
[Unit]
Description=Foo
[Swap]
SyslogLevel=notice

