# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_11.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=11

