# tests/auto/InvalidValue/2.31/good/Swap_LimitFSIZE_0.target
[Unit]
Description=Foo
[Swap]
LimitFSIZE=0

