# tests/auto/InvalidValue/2.31/good/Swap_ProtectHome_read-only.target
[Unit]
Description=Foo
[Swap]
ProtectHome=read-only

