# tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_6h.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=6h

