# tests/auto/InvalidValue/2.31/good/Swap_LimitRSS_infinity.target
[Unit]
Description=Foo
[Swap]
LimitRSS=infinity

