# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_12.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=12

