# tests/auto/InvalidValue/2.31/good/Swap_BlockIODeviceWeight_--dev--foo%2010000.target
[Unit]
Description=Foo
[Swap]
BlockIODeviceWeight=/dev/foo 10000

