# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENFILE.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENFILE

