# tests/auto/InvalidValue/2.31/good/Swap_MemoryAccounting_0.target
[Unit]
Description=Foo
[Swap]
MemoryAccounting=0

