# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EEXIST.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EEXIST

