# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGPWR.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGPWR

