# tests/auto/InvalidValue/2.31/good/Swap_TTYReset_on.target
[Unit]
Description=Foo
[Swap]
TTYReset=on

