# tests/auto/InvalidValue/2.31/good/Swap_StandardOutput_abc.target
[Unit]
Description=Foo
[Swap]
StandardOutput=abc

