# tests/auto/InvalidValue/2.31/good/Swap_SyslogFacility_syslog.target
[Unit]
Description=Foo
[Swap]
SyslogFacility=syslog

