# tests/auto/InvalidValue/2.31/good/Swap_LimitDATA_infinity.target
[Unit]
Description=Foo
[Swap]
LimitDATA=infinity

