# tests/auto/InvalidValue/2.31/good/Swap_MemoryDenyWriteExecute_1.target
[Unit]
Description=Foo
[Swap]
MemoryDenyWriteExecute=1

