# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGTSTP.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGTSTP

