# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGUSR1.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGUSR1

