# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGIOT.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGIOT

