# tests/auto/InvalidValue/2.31/good/Swap_CPUAffinity_numa.target
[Unit]
Description=Foo
[Swap]
CPUAffinity=numa

