# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGBUS.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGBUS

