# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EMULTIHOP.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EMULTIHOP

