# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EBUSY.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EBUSY

