# tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGTTOU.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGTTOU

