# tests/auto/InvalidValue/2.31/good/Swap_RestrictSUIDSGID_no.target
[Unit]
Description=Foo
[Swap]
RestrictSUIDSGID=no

