# tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGXFSZ.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGXFSZ

