# tests/auto/InvalidValue/2.31/good/Swap_CPUSchedulingPolicy_rr.target
[Unit]
Description=Foo
[Swap]
CPUSchedulingPolicy=rr

