# tests/auto/InvalidValue/2.31/good/Swap_KillSignal_47.target
[Unit]
Description=Foo
[Swap]
KillSignal=47

