# tests/auto/InvalidValue/2.31/good/Swap_BlockIOAccounting_1.target
[Unit]
Description=Foo
[Swap]
BlockIOAccounting=1

