# tests/auto/InvalidValue/2.31/good/Swap_DevicePolicy_closed.target
[Unit]
Description=Foo
[Swap]
DevicePolicy=closed

