# tests/auto/InvalidValue/2.31/good/Swap_StartupCPUWeight_10000.target
[Unit]
Description=Foo
[Swap]
StartupCPUWeight=10000

