# tests/auto/InvalidValue/2.31/good/Swap_MemoryHigh_9999999999999M.target
[Unit]
Description=Foo
[Swap]
MemoryHigh=9999999999999M

