# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOTRECOVERABLE.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOTRECOVERABLE

