# tests/auto/InvalidValue/2.31/good/Swap_MemorySwapMax_0B.target
[Unit]
Description=Foo
[Swap]
MemorySwapMax=0B

