# tests/auto/InvalidValue/2.31/good/Swap_DeviceAllow_--dev--bar.target
[Unit]
Description=Foo
[Swap]
DeviceAllow=/dev/bar

