# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EBADR.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EBADR

