# tests/auto/InvalidValue/2.31/good/Swap_SendSIGHUP_yes.target
[Unit]
Description=Foo
[Swap]
SendSIGHUP=yes

