# tests/auto/InvalidValue/2.31/good/Swap_BlockIOWriteBandwidthMax_--dev--bar%209999999999999K.target
[Unit]
Description=Foo
[Swap]
BlockIOWriteBandwidthMax=/dev/bar 9999999999999K

