# tests/auto/InvalidValue/2.31/good/Swap_Environment_a%3Db.target
[Unit]
Description=Foo
[Swap]
Environment=a=b

