# tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_8minute.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=8minute

