# tests/auto/InvalidValue/2.31/good/Swap_InaccessableDirectories_abc.target
[Unit]
Description=Foo
[Swap]
InaccessableDirectories=abc

