# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGSYS.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGSYS

