# tests/auto/InvalidValue/2.31/good/Swap_IOWeight_1.target
[Unit]
Description=Foo
[Swap]
IOWeight=1

