# tests/auto/InvalidValue/2.31/good/Swap_RemoveIPC_1.target
[Unit]
Description=Foo
[Swap]
RemoveIPC=1

