# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EHOSTUNREACH.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EHOSTUNREACH

