# tests/auto/InvalidValue/2.31/good/Swap_KillSignal_24.target
[Unit]
Description=Foo
[Swap]
KillSignal=24

