# tests/auto/InvalidValue/2.31/good/Swap_MemoryLimit_9999999999999T.target
[Unit]
Description=Foo
[Swap]
MemoryLimit=9999999999999T

