# tests/auto/InvalidValue/2.31/good/Swap_IOSchedulingClass_realtime.target
[Unit]
Description=Foo
[Swap]
IOSchedulingClass=realtime

