# tests/auto/InvalidValue/2.31/good/Swap_RemoveIPC_true.target
[Unit]
Description=Foo
[Swap]
RemoveIPC=true

