# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_31.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=31

