# tests/auto/InvalidValue/2.31/good/Swap_TTYReset_yes.target
[Unit]
Description=Foo
[Swap]
TTYReset=yes

