# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGLOST.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGLOST

