# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EKEYREVOKED.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EKEYREVOKED

