# tests/auto/InvalidValue/2.31/good/Swap_LimitMEMLOCK_infinity.target
[Unit]
Description=Foo
[Swap]
LimitMEMLOCK=infinity

