# tests/auto/InvalidValue/2.31/good/Swap_KillSignal_14.target
[Unit]
Description=Foo
[Swap]
KillSignal=14

