# tests/auto/InvalidValue/2.31/good/Swap_ProtectHome_on.target
[Unit]
Description=Foo
[Swap]
ProtectHome=on

