# tests/auto/InvalidValue/2.31/good/Swap_StandardError_null.target
[Unit]
Description=Foo
[Swap]
StandardError=null

