# tests/auto/InvalidValue/2.31/good/Swap_LimitLOCKS_0.target
[Unit]
Description=Foo
[Swap]
LimitLOCKS=0

