# tests/auto/InvalidValue/2.31/good/Swap_LimitSIGPENDING_infinity.target
[Unit]
Description=Foo
[Swap]
LimitSIGPENDING=infinity

