# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGINT.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGINT

