# tests/auto/InvalidValue/2.31/good/Swap_RestrictSUIDSGID_true.target
[Unit]
Description=Foo
[Swap]
RestrictSUIDSGID=true

