# tests/auto/InvalidValue/2.31/good/Swap_IOAccounting_on.target
[Unit]
Description=Foo
[Swap]
IOAccounting=on

