# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENETUNREACH.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENETUNREACH

