# tests/auto/InvalidValue/2.31/good/Swap_Group_metoo.target
[Unit]
Description=Foo
[Swap]
Group=metoo

