# tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGPROF.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGPROF

