# tests/auto/InvalidValue/2.31/good/Swap_TTYReset_0.target
[Unit]
Description=Foo
[Swap]
TTYReset=0

