# tests/auto/InvalidValue/2.31/good/Swap_LimitSIGPENDING_9999999999999.target
[Unit]
Description=Foo
[Swap]
LimitSIGPENDING=9999999999999

