# tests/auto/InvalidValue/2.31/good/Swap_Personality_ppc64-le.target
[Unit]
Description=Foo
[Swap]
Personality=ppc64-le

