# tests/auto/InvalidValue/2.31/good/Swap_LimitNPROC_infinity.target
[Unit]
Description=Foo
[Swap]
LimitNPROC=infinity

