# tests/auto/InvalidValue/2.31/good/Swap_LimitRTTIME_9999999999999.target
[Unit]
Description=Foo
[Swap]
LimitRTTIME=9999999999999

