# tests/auto/InvalidValue/2.31/good/Swap_TTYVTDisallocate_on.target
[Unit]
Description=Foo
[Swap]
TTYVTDisallocate=on

