# tests/auto/InvalidValue/2.31/good/Swap_RemoveIPC_on.target
[Unit]
Description=Foo
[Swap]
RemoveIPC=on

