# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EROFS.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EROFS

