# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGWINCH.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGWINCH

