# tests/auto/InvalidValue/2.31/good/Swap_IOWeight_10000.target
[Unit]
Description=Foo
[Swap]
IOWeight=10000

