# tests/auto/InvalidValue/2.31/good/Swap_LimitMSGQUEUE_9999999999999.target
[Unit]
Description=Foo
[Swap]
LimitMSGQUEUE=9999999999999

