# tests/auto/InvalidValue/2.31/good/Swap_LimitNOFILE_0.target
[Unit]
Description=Foo
[Swap]
LimitNOFILE=0

