# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ECHRNG.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ECHRNG

