# tests/auto/InvalidValue/2.31/good/Swap_SendSIGHUP_no.target
[Unit]
Description=Foo
[Swap]
SendSIGHUP=no

