# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGQUIT.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGQUIT

