# tests/auto/InvalidValue/2.31/good/Swap_User_metoo.target
[Unit]
Description=Foo
[Swap]
User=metoo

