# tests/auto/InvalidValue/2.31/good/Swap_RestrictRealtime_off.target
[Unit]
Description=Foo
[Swap]
RestrictRealtime=off

