# tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_16seconds.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=16seconds

