# tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_24w.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=24w

