# tests/auto/InvalidValue/2.31/good/Swap_Delegate_pids.target
[Unit]
Description=Foo
[Swap]
Delegate=pids

