# tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_infinity.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=infinity

