# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ERESTART.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ERESTART

