# tests/auto/InvalidValue/2.31/good/Swap_PrivateDevices_off.target
[Unit]
Description=Foo
[Swap]
PrivateDevices=off

