# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGCONT.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGCONT

