# tests/auto/InvalidValue/2.31/good/Swap_MemoryAccounting_true.target
[Unit]
Description=Foo
[Swap]
MemoryAccounting=true

