# tests/auto/InvalidValue/2.31/good/Swap_RootDirectory_--usr--sbin.target
[Unit]
Description=Foo
[Swap]
RootDirectory=/usr/sbin

