# tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGUNUSED.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGUNUSED

