# tests/auto/InvalidValue/2.31/good/Swap_TTYReset_1.target
[Unit]
Description=Foo
[Swap]
TTYReset=1

