# tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_12month.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=12month

