# tests/auto/InvalidValue/2.31/good/Swap_MemoryAccounting_1.target
[Unit]
Description=Foo
[Swap]
MemoryAccounting=1

