# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGVTALRM.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGVTALRM

