# tests/auto/InvalidValue/2.31/good/Swap_SyslogFacility_ftp.target
[Unit]
Description=Foo
[Swap]
SyslogFacility=ftp

