# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOANO.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOANO

