# tests/auto/InvalidValue/2.31/good/Swap_KillMode_none.target
[Unit]
Description=Foo
[Swap]
KillMode=none

