# tests/auto/InvalidValue/2.31/good/Swap_IODeviceLatencyTargetSec_abc.target
[Unit]
Description=Foo
[Swap]
IODeviceLatencyTargetSec=abc

