# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EXFULL.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EXFULL

