# tests/auto/InvalidValue/2.31/good/Swap_DisableControllers_io.target
[Unit]
Description=Foo
[Swap]
DisableControllers=io

