# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGTTIN.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGTTIN

