# tests/auto/InvalidValue/2.31/good/Swap_RestrictRealtime_0.target
[Unit]
Description=Foo
[Swap]
RestrictRealtime=0

