# tests/auto/InvalidValue/2.31/good/Swap_LimitCPU_infinity.target
[Unit]
Description=Foo
[Swap]
LimitCPU=infinity

