# tests/auto/InvalidValue/2.31/good/Swap_SupplementaryGroups_iamagroup.target
[Unit]
Description=Foo
[Swap]
SupplementaryGroups=iamagroup

