# tests/auto/InvalidValue/2.31/good/Swap_IgnoreSIGPIPE_off.target
[Unit]
Description=Foo
[Swap]
IgnoreSIGPIPE=off

