# tests/auto/InvalidValue/2.31/good/Swap_SendSIGHUP_0.target
[Unit]
Description=Foo
[Swap]
SendSIGHUP=0

