# tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_19s.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=19s

