# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_13.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=13

