# tests/auto/InvalidValue/2.31/good/Swap_CPUShares_262144.target
[Unit]
Description=Foo
[Swap]
CPUShares=262144

