# tests/auto/InvalidValue/2.31/good/Swap_PrivateDevices_1.target
[Unit]
Description=Foo
[Swap]
PrivateDevices=1

