# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENETRESET.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENETRESET

