# tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGVTALRM.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGVTALRM

