# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGPWR.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGPWR

