# tests/auto/InvalidValue/2.31/good/Swap_DisableControllers_devices.target
[Unit]
Description=Foo
[Swap]
DisableControllers=devices

