# tests/auto/InvalidValue/2.31/good/Swap_SyslogLevel_warning.target
[Unit]
Description=Foo
[Swap]
SyslogLevel=warning

