# tests/auto/InvalidValue/2.31/good/Swap_LimitAS_infinity.target
[Unit]
Description=Foo
[Swap]
LimitAS=infinity

