# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGURG.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGURG

