# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGPIPE.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGPIPE

