# tests/auto/InvalidValue/2.31/good/Swap_LimitSTACK_9999999999999.target
[Unit]
Description=Foo
[Swap]
LimitSTACK=9999999999999

