# tests/auto/InvalidValue/2.31/good/Swap_SyslogLevelPrefix_yes.target
[Unit]
Description=Foo
[Swap]
SyslogLevelPrefix=yes

