# tests/auto/InvalidValue/2.31/good/Swap_SendSIGHUP_off.target
[Unit]
Description=Foo
[Swap]
SendSIGHUP=off

