# tests/auto/InvalidValue/2.31/good/Swap_StandardOutput_tty.target
[Unit]
Description=Foo
[Swap]
StandardOutput=tty

