# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EBADMSG.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EBADMSG

