# tests/auto/InvalidValue/2.31/good/Swap_PrivateNetwork_0.target
[Unit]
Description=Foo
[Swap]
PrivateNetwork=0

