# tests/auto/InvalidValue/2.31/good/Swap_SyslogLevelPrefix_true.target
[Unit]
Description=Foo
[Swap]
SyslogLevelPrefix=true

