# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGLOST.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGLOST

