# tests/auto/InvalidValue/2.31/good/Swap_SyslogFacility_user.target
[Unit]
Description=Foo
[Swap]
SyslogFacility=user

