# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EPROTONOSUPPORT.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EPROTONOSUPPORT

