# tests/auto/InvalidValue/2.31/good/Swap_TasksAccounting_true.target
[Unit]
Description=Foo
[Swap]
TasksAccounting=true

