# tests/auto/InvalidValue/2.31/good/Swap_RestrictRealtime_no.target
[Unit]
Description=Foo
[Swap]
RestrictRealtime=no

