# tests/auto/InvalidValue/2.31/good/Swap_Group_iamagroup.target
[Unit]
Description=Foo
[Swap]
Group=iamagroup

