# tests/auto/InvalidValue/2.31/good/Swap_SendSIGHUP_on.target
[Unit]
Description=Foo
[Swap]
SendSIGHUP=on

