# tests/auto/InvalidValue/2.31/good/Swap_CPUSchedulingPolicy_fifo.target
[Unit]
Description=Foo
[Swap]
CPUSchedulingPolicy=fifo

