# tests/auto/InvalidValue/2.31/good/Swap_Personality_s390.target
[Unit]
Description=Foo
[Swap]
Personality=s390

