# tests/auto/InvalidValue/2.31/good/Swap_MemoryDenyWriteExecute_off.target
[Unit]
Description=Foo
[Swap]
MemoryDenyWriteExecute=off

