# tests/auto/InvalidValue/2.31/good/Swap_ProtectSystem_1.target
[Unit]
Description=Foo
[Swap]
ProtectSystem=1

