# tests/auto/InvalidValue/2.31/good/Swap_LimitRTPRIO_9999999999999.target
[Unit]
Description=Foo
[Swap]
LimitRTPRIO=9999999999999

