# tests/auto/InvalidValue/2.31/good/Swap_RestrictAddressFamilies_abc.target
[Unit]
Description=Foo
[Swap]
RestrictAddressFamilies=abc

