# tests/auto/InvalidValue/2.31/good/Swap_BlockIOAccounting_true.target
[Unit]
Description=Foo
[Swap]
BlockIOAccounting=true

