# tests/auto/InvalidValue/2.31/good/Swap_SystemCallFilter_abc.target
[Unit]
Description=Foo
[Swap]
SystemCallFilter=abc

