# tests/auto/InvalidValue/2.31/good/Swap_PrivateTmp_on.target
[Unit]
Description=Foo
[Swap]
PrivateTmp=on

