# tests/auto/InvalidValue/2.31/good/Swap_MemoryDenyWriteExecute_true.target
[Unit]
Description=Foo
[Swap]
MemoryDenyWriteExecute=true

