# tests/auto/InvalidValue/2.31/good/Swap_TasksAccounting_yes.target
[Unit]
Description=Foo
[Swap]
TasksAccounting=yes

