# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EPIPE.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EPIPE

