# tests/auto/InvalidValue/2.31/good/Swap_BlockIOReadBandwidthMax_--dev--foo%200K.target
[Unit]
Description=Foo
[Swap]
BlockIOReadBandwidthMax=/dev/foo 0K

