# tests/auto/InvalidValue/2.31/good/Swap_BlockIOAccounting_no.target
[Unit]
Description=Foo
[Swap]
BlockIOAccounting=no

