# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EBADFD.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EBADFD

