# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ELNRANGE.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ELNRANGE

