# tests/auto/InvalidValue/2.31/good/Swap_SyslogLevelPrefix_false.target
[Unit]
Description=Foo
[Swap]
SyslogLevelPrefix=false

