# tests/auto/InvalidValue/2.31/good/Swap_BlockIOAccounting_off.target
[Unit]
Description=Foo
[Swap]
BlockIOAccounting=off

