# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGUNUSED.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGUNUSED

