# tests/auto/InvalidValue/2.31/good/Swap_CPUQuota_9999999999999%25.target
[Unit]
Description=Foo
[Swap]
CPUQuota=9999999999999%

