# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EILSEQ.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EILSEQ

