# tests/auto/InvalidValue/2.31/good/Swap_Delegate_devices.target
[Unit]
Description=Foo
[Swap]
Delegate=devices

