# tests/auto/InvalidValue/2.31/good/Swap_StandardError_journal.target
[Unit]
Description=Foo
[Swap]
StandardError=journal

