# tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGBUS.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGBUS

