# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGKILL.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGKILL

