# tests/auto/InvalidValue/2.31/good/Swap_BlockIOWeight_1.target
[Unit]
Description=Foo
[Swap]
BlockIOWeight=1

