# tests/auto/InvalidValue/2.31/good/Swap_BlockIOReadBandwidthMax_--dev--foo%209999999999999G.target
[Unit]
Description=Foo
[Swap]
BlockIOReadBandwidthMax=/dev/foo 9999999999999G

