# tests/auto/InvalidValue/2.31/good/Swap_RestrictSUIDSGID_0.target
[Unit]
Description=Foo
[Swap]
RestrictSUIDSGID=0

