# tests/auto/InvalidValue/2.31/good/Swap_UtmpMode_login.target
[Unit]
Description=Foo
[Swap]
UtmpMode=login

