# tests/auto/InvalidValue/2.31/good/Swap_PrivateNetwork_yes.target
[Unit]
Description=Foo
[Swap]
PrivateNetwork=yes

