# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGCHLD.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGCHLD

