# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGALRM.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGALRM

