# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EOWNERDEAD.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EOWNERDEAD

