# tests/auto/InvalidValue/2.31/good/Swap_SecureBits_keep-caps.target
[Unit]
Description=Foo
[Swap]
SecureBits=keep-caps

