# tests/auto/InvalidValue/2.31/good/Swap_IgnoreSIGPIPE_on.target
[Unit]
Description=Foo
[Swap]
IgnoreSIGPIPE=on

