# tests/auto/InvalidValue/2.31/good/Swap_RemoveIPC_no.target
[Unit]
Description=Foo
[Swap]
RemoveIPC=no

