# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGTTOU.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGTTOU

