# tests/auto/InvalidValue/2.31/good/Swap_CPUSchedulingPolicy_other.target
[Unit]
Description=Foo
[Swap]
CPUSchedulingPolicy=other

