# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGINFO.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGINFO

