# tests/auto/InvalidValue/2.31/good/Swap_MountAPIVFS_1.target
[Unit]
Description=Foo
[Swap]
MountAPIVFS=1

