# tests/auto/InvalidValue/2.31/good/Swap_DisableControllers_blkio.target
[Unit]
Description=Foo
[Swap]
DisableControllers=blkio

