# tests/auto/InvalidValue/2.31/good/Swap_SyslogFacility_auth.target
[Unit]
Description=Foo
[Swap]
SyslogFacility=auth

