# tests/auto/InvalidValue/2.31/good/Swap_MountAPIVFS_off.target
[Unit]
Description=Foo
[Swap]
MountAPIVFS=off

