# tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGIOT.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGIOT

