# tests/auto/InvalidValue/2.31/good/Swap_IOAccounting_0.target
[Unit]
Description=Foo
[Swap]
IOAccounting=0

