# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGBUS.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGBUS

