# tests/auto/InvalidValue/2.31/good/Swap_DisableControllers_memory.target
[Unit]
Description=Foo
[Swap]
DisableControllers=memory

