# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOPKG.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOPKG

