# tests/auto/InvalidValue/2.31/good/Swap_TTYVHangup_true.target
[Unit]
Description=Foo
[Swap]
TTYVHangup=true

