# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ERFKILL.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ERFKILL

