# tests/auto/InvalidValue/2.31/good/Swap_IOAccounting_true.target
[Unit]
Description=Foo
[Swap]
IOAccounting=true

