# tests/auto/InvalidValue/2.31/good/Swap_RestrictSUIDSGID_off.target
[Unit]
Description=Foo
[Swap]
RestrictSUIDSGID=off

