# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EFBIG.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EFBIG

