# tests/auto/InvalidValue/2.31/good/Swap_MemoryLow_0K.target
[Unit]
Description=Foo
[Swap]
MemoryLow=0K

