# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EL2NSYNC.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EL2NSYNC

