# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EMLINK.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EMLINK

