# tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_4hour.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=4hour

