# tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_10M.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=10M

