# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ESOCKTNOSUPPORT.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ESOCKTNOSUPPORT

