# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ECOMM.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ECOMM

