# tests/auto/InvalidValue/2.31/good/Swap_TimeoutSec_25years.target
[Unit]
Description=Foo
[Swap]
TimeoutSec=25years

