# tests/auto/InvalidValue/2.31/good/Swap_IOWriteIOPSMax_abc.target
[Unit]
Description=Foo
[Swap]
IOWriteIOPSMax=abc

