# tests/auto/InvalidValue/2.31/good/Swap_SendSIGKILL_off.target
[Unit]
Description=Foo
[Swap]
SendSIGKILL=off

