# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGUSR2.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGUSR2

