# tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGUSR1.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGUSR1

