# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGSTKFLT.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGSTKFLT

