# tests/auto/InvalidValue/2.31/good/Swap_PrivateNetwork_1.target
[Unit]
Description=Foo
[Swap]
PrivateNetwork=1

