# tests/auto/InvalidValue/2.31/good/Swap_CPUSchedulingPolicy_idle.target
[Unit]
Description=Foo
[Swap]
CPUSchedulingPolicy=idle

