# tests/auto/InvalidValue/2.31/good/Swap_IOSchedulingClass_best-effort.target
[Unit]
Description=Foo
[Swap]
IOSchedulingClass=best-effort

