# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_SIGUNUSED.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=SIGUNUSED

