# tests/auto/InvalidValue/2.31/good/Swap_SendSIGKILL_0.target
[Unit]
Description=Foo
[Swap]
SendSIGKILL=0

