# tests/auto/InvalidValue/2.31/good/Swap_ProtectHome_yes.target
[Unit]
Description=Foo
[Swap]
ProtectHome=yes

