# tests/auto/InvalidValue/2.31/good/Swap_MountAPIVFS_on.target
[Unit]
Description=Foo
[Swap]
MountAPIVFS=on

