# tests/auto/InvalidValue/2.31/good/Swap_LimitAS_0.target
[Unit]
Description=Foo
[Swap]
LimitAS=0

