# tests/auto/InvalidValue/2.31/good/Swap_ProtectHome_false.target
[Unit]
Description=Foo
[Swap]
ProtectHome=false

