# tests/auto/InvalidValue/2.31/good/Swap_TasksAccounting_on.target
[Unit]
Description=Foo
[Swap]
TasksAccounting=on

