# tests/auto/InvalidValue/2.31/good/Swap_TasksMax_9999999999999.target
[Unit]
Description=Foo
[Swap]
TasksMax=9999999999999

