# tests/auto/InvalidValue/2.31/good/Swap_TTYReset_true.target
[Unit]
Description=Foo
[Swap]
TTYReset=true

