# tests/auto/InvalidValue/2.31/good/Swap_SyslogFacility_mail.target
[Unit]
Description=Foo
[Swap]
SyslogFacility=mail

