# tests/auto/InvalidValue/2.31/good/Swap_WatchdogSignal_22.target
[Unit]
Description=Foo
[Swap]
WatchdogSignal=22

