# tests/auto/InvalidValue/2.31/good/Swap_PassEnvironment_abc.target
[Unit]
Description=Foo
[Swap]
PassEnvironment=abc

