# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_48.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=48

