# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EPERM.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EPERM

