# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOTSUP.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOTSUP

