# tests/auto/InvalidValue/2.31/good/Swap_Options_abc.target
[Unit]
Description=Foo
[Swap]
Options=abc

