# tests/auto/InvalidValue/2.31/good/Swap_LimitNICE_9999999999999.target
[Unit]
Description=Foo
[Swap]
LimitNICE=9999999999999

