# tests/auto/InvalidValue/2.31/good/Swap_MemoryAccounting_no.target
[Unit]
Description=Foo
[Swap]
MemoryAccounting=no

