# tests/auto/InvalidValue/2.31/good/Swap_MemoryAccounting_false.target
[Unit]
Description=Foo
[Swap]
MemoryAccounting=false

