# tests/auto/InvalidValue/2.31/good/Swap_LimitFSIZE_infinity.target
[Unit]
Description=Foo
[Swap]
LimitFSIZE=infinity

