# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_EBADF.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=EBADF

