# tests/auto/InvalidValue/2.31/good/Swap_TasksAccounting_false.target
[Unit]
Description=Foo
[Swap]
TasksAccounting=false

