# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOBUFS.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOBUFS

