# tests/auto/InvalidValue/2.31/good/Swap_StandardOutput_kmsg%2Bconsole.target
[Unit]
Description=Foo
[Swap]
StandardOutput=kmsg+console

