# tests/auto/InvalidValue/2.31/good/Swap_BlockIOReadBandwidthMax_--dev--bar%200G.target
[Unit]
Description=Foo
[Swap]
BlockIOReadBandwidthMax=/dev/bar 0G

