# tests/auto/InvalidValue/2.31/good/Swap_ProtectSystem_no.target
[Unit]
Description=Foo
[Swap]
ProtectSystem=no

