# tests/auto/InvalidValue/2.31/good/Swap_TasksAccounting_0.target
[Unit]
Description=Foo
[Swap]
TasksAccounting=0

