# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGXCPU.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGXCPU

