# tests/auto/InvalidValue/2.31/good/Swap_FinalKillSignal_SIGABRT.target
[Unit]
Description=Foo
[Swap]
FinalKillSignal=SIGABRT

