# tests/auto/InvalidValue/2.31/good/Swap_SyslogLevelPrefix_off.target
[Unit]
Description=Foo
[Swap]
SyslogLevelPrefix=off

