# tests/auto/InvalidValue/2.31/good/Swap_KillMode_control-group.target
[Unit]
Description=Foo
[Swap]
KillMode=control-group

