# tests/auto/InvalidValue/2.31/good/Swap_InaccessiblePaths_abc.target
[Unit]
Description=Foo
[Swap]
InaccessiblePaths=abc

