# tests/auto/InvalidValue/2.31/good/Swap_TTYReset_false.target
[Unit]
Description=Foo
[Swap]
TTYReset=false

