# tests/auto/InvalidValue/2.31/good/Swap_KillSignal_SIGTSTP.target
[Unit]
Description=Foo
[Swap]
KillSignal=SIGTSTP

