# tests/auto/InvalidValue/2.31/good/Swap_SystemCallErrorNumber_ENOKEY.target
[Unit]
Description=Foo
[Swap]
SystemCallErrorNumber=ENOKEY

