# tests/auto/InvalidValue/2.31/good/Tap_VNetHeader_false.target
[Unit]
Description=Foo
[Tap]
VNetHeader=false

