# tests/auto/InvalidValue/2.31/good/Tap_PacketInfo_no.target
[Unit]
Description=Foo
[Tap]
PacketInfo=no

