# tests/auto/InvalidValue/2.31/good/Tap_MultiQueue_0.target
[Unit]
Description=Foo
[Tap]
MultiQueue=0

