# tests/auto/InvalidValue/2.31/good/Tap_MultiQueue_yes.target
[Unit]
Description=Foo
[Tap]
MultiQueue=yes

