# tests/auto/InvalidValue/2.31/good/Tap_VNetHeader_true.target
[Unit]
Description=Foo
[Tap]
VNetHeader=true

