# tests/auto/InvalidValue/2.31/good/Tap_OneQueue_off.target
[Unit]
Description=Foo
[Tap]
OneQueue=off

