# tests/auto/InvalidValue/2.31/good/Tap_User_iamauser.target
[Unit]
Description=Foo
[Tap]
User=iamauser

