# tests/auto/InvalidValue/2.31/good/Tap_MultiQueue_1.target
[Unit]
Description=Foo
[Tap]
MultiQueue=1

