# tests/auto/InvalidValue/2.31/good/Tap_MultiQueue_on.target
[Unit]
Description=Foo
[Tap]
MultiQueue=on

