# tests/auto/InvalidValue/2.31/good/Tap_PacketInfo_true.target
[Unit]
Description=Foo
[Tap]
PacketInfo=true

