# tests/auto/InvalidValue/2.31/good/Tap_Group_metoo.target
[Unit]
Description=Foo
[Tap]
Group=metoo

