# tests/auto/InvalidValue/2.31/good/Tap_PacketInfo_off.target
[Unit]
Description=Foo
[Tap]
PacketInfo=off

