# tests/auto/InvalidValue/2.31/good/Tap_VNetHeader_yes.target
[Unit]
Description=Foo
[Tap]
VNetHeader=yes

