# tests/auto/InvalidValue/2.31/good/Tap_User_metoo.target
[Unit]
Description=Foo
[Tap]
User=metoo

