# tests/auto/InvalidValue/2.31/good/Time_RootDistanceMaxSec_0.target
[Unit]
Description=Foo
[Time]
RootDistanceMaxSec=0

