# tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_0days.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=0days

