# tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_26year.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=26year

