# tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_6h.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=6h

