# tests/auto/InvalidValue/2.31/good/Timer_OnUnitActiveSec_17second.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=17second

