# tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_13months.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=13months

