# tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_2d.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=2d

