# tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_28%C2%B5s.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=28µs

