# tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_1day.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=1day

