# tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_9min.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=9min

