# tests/auto/InvalidValue/2.31/good/Timer_OnCalendar_abc.target
[Unit]
Description=Foo
[Timer]
OnCalendar=abc

