# tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_18sec.target
[Unit]
Description=Foo
[Timer]
OnBootSec=18sec

