# tests/auto/InvalidValue/2.31/good/Timer_OnClockChange_yes.target
[Unit]
Description=Foo
[Timer]
OnClockChange=yes

