# tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_14msec.target
[Unit]
Description=Foo
[Timer]
AccuracySec=14msec

