# tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_11m.target
[Unit]
Description=Foo
[Timer]
AccuracySec=11m

