# tests/auto/InvalidValue/2.31/good/Timer_RemainAfterElapse_no.target
[Unit]
Description=Foo
[Timer]
RemainAfterElapse=no

