# tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_15ms.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=15ms

