# tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_4hour.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=4hour

