# tests/auto/InvalidValue/2.31/good/Timer_OnUnitActiveSec_22weeks.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=22weeks

