# tests/auto/InvalidValue/2.31/good/Timer_OnUnitActiveSec_25years.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=25years

