# tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_23week.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=23week

