# tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_0days.target
[Unit]
Description=Foo
[Timer]
AccuracySec=0days

