# tests/auto/InvalidValue/2.31/good/Timer_Persistent_0.target
[Unit]
Description=Foo
[Timer]
Persistent=0

