# tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_infinity.target
[Unit]
Description=Foo
[Timer]
OnBootSec=infinity

