# tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_16seconds.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=16seconds

