# tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_15ms.target
[Unit]
Description=Foo
[Timer]
AccuracySec=15ms

