# tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_17second.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=17second

