# tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_12month.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=12month

