# tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_5hr.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=5hr

