# tests/auto/InvalidValue/2.31/good/Time_RootDistanceMaxSec_9999999999999.target
[Unit]
Description=Foo
[Time]
RootDistanceMaxSec=9999999999999

