# tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_11m.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=11m

