# tests/auto/InvalidValue/2.31/good/Timer_OnTimezoneChange_true.target
[Unit]
Description=Foo
[Timer]
OnTimezoneChange=true

