# tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_8minute.target
[Unit]
Description=Foo
[Timer]
AccuracySec=8minute

