# tests/auto/InvalidValue/2.31/good/Timer_RemainAfterElapse_yes.target
[Unit]
Description=Foo
[Timer]
RemainAfterElapse=yes

