# tests/auto/InvalidValue/2.31/good/Timer_OnTimezoneChange_false.target
[Unit]
Description=Foo
[Timer]
OnTimezoneChange=false

