# tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_16seconds.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=16seconds

