# tests/auto/InvalidValue/2.31/good/Timer_OnUnitActiveSec_28%C2%B5s.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=28µs

