# tests/auto/InvalidValue/2.31/good/Timer_WakeSystem_yes.target
[Unit]
Description=Foo
[Timer]
WakeSystem=yes

