# tests/auto/InvalidValue/2.31/good/Timer_OnUnitActiveSec_15ms.target
[Unit]
Description=Foo
[Timer]
OnUnitActiveSec=15ms

