# tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_8minute.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=8minute

