# tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_27y.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=27y

