# tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_7minutes.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=7minutes

