# tests/auto/InvalidValue/2.31/good/Timer_Persistent_no.target
[Unit]
Description=Foo
[Timer]
Persistent=no

