# tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_13months.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=13months

