# tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_infinity.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=infinity

