# tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_21us.target
[Unit]
Description=Foo
[Timer]
OnBootSec=21us

