# tests/auto/InvalidValue/2.31/good/Timer_OnClockChange_off.target
[Unit]
Description=Foo
[Timer]
OnClockChange=off

