# tests/auto/InvalidValue/2.31/good/Timer_Persistent_1.target
[Unit]
Description=Foo
[Timer]
Persistent=1

