# tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_15ms.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=15ms

