# tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_9min.target
[Unit]
Description=Foo
[Timer]
AccuracySec=9min

