# tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_15ms.target
[Unit]
Description=Foo
[Timer]
OnBootSec=15ms

