# tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_10M.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=10M

