# tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_18sec.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=18sec

