# tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_14msec.target
[Unit]
Description=Foo
[Timer]
OnBootSec=14msec

