# tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_4hour.target
[Unit]
Description=Foo
[Timer]
AccuracySec=4hour

