# tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_12month.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=12month

