# tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_21us.target
[Unit]
Description=Foo
[Timer]
AccuracySec=21us

