# tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_3hours.target
[Unit]
Description=Foo
[Timer]
OnBootSec=3hours

