# tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_22weeks.target
[Unit]
Description=Foo
[Timer]
AccuracySec=22weeks

