# tests/auto/InvalidValue/2.31/good/Timer_Persistent_false.target
[Unit]
Description=Foo
[Timer]
Persistent=false

