# tests/auto/InvalidValue/2.31/good/Timer_RemainAfterElapse_true.target
[Unit]
Description=Foo
[Timer]
RemainAfterElapse=true

