# tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_20usec.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=20usec

