# tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_1day.target
[Unit]
Description=Foo
[Timer]
AccuracySec=1day

