# tests/auto/InvalidValue/2.31/good/Timer_RemainAfterElapse_off.target
[Unit]
Description=Foo
[Timer]
RemainAfterElapse=off

