# tests/auto/InvalidValue/2.31/good/Timer_WakeSystem_off.target
[Unit]
Description=Foo
[Timer]
WakeSystem=off

