# tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_9min.target
[Unit]
Description=Foo
[Timer]
OnBootSec=9min

