# tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_18sec.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=18sec

