# tests/auto/InvalidValue/2.31/good/Timer_AccuracySec_20usec.target
[Unit]
Description=Foo
[Timer]
AccuracySec=20usec

