# tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_27y.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=27y

