# tests/auto/InvalidValue/2.31/good/Timer_RandomizedDelaySec_12month.target
[Unit]
Description=Foo
[Timer]
RandomizedDelaySec=12month

