# tests/auto/InvalidValue/2.31/good/Timer_OnTimezoneChange_yes.target
[Unit]
Description=Foo
[Timer]
OnTimezoneChange=yes

