# tests/auto/InvalidValue/2.31/good/Time_FallbackNTP_abc.target
[Unit]
Description=Foo
[Time]
FallbackNTP=abc

