# tests/auto/InvalidValue/2.31/good/Timer_OnActiveSec_22weeks.target
[Unit]
Description=Foo
[Timer]
OnActiveSec=22weeks

