# tests/auto/InvalidValue/2.31/good/Timer_OnTimezoneChange_no.target
[Unit]
Description=Foo
[Timer]
OnTimezoneChange=no

