# tests/auto/InvalidValue/2.31/good/Timer_OnUnitInactiveSec_14msec.target
[Unit]
Description=Foo
[Timer]
OnUnitInactiveSec=14msec

