# tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_0days.target
[Unit]
Description=Foo
[Timer]
OnBootSec=0days

