# tests/auto/InvalidValue/2.31/good/Timer_OnBootSec_12month.target
[Unit]
Description=Foo
[Timer]
OnBootSec=12month

