# tests/auto/InvalidValue/2.31/good/Tun_MultiQueue_on.target
[Unit]
Description=Foo
[Tun]
MultiQueue=on

