# tests/auto/InvalidValue/2.31/good/Tun_PacketInfo_yes.target
[Unit]
Description=Foo
[Tun]
PacketInfo=yes

