# tests/auto/InvalidValue/2.31/good/Tun_MultiQueue_false.target
[Unit]
Description=Foo
[Tun]
MultiQueue=false

