# tests/auto/InvalidValue/2.31/good/Tun_PacketInfo_false.target
[Unit]
Description=Foo
[Tun]
PacketInfo=false

