# tests/auto/InvalidValue/2.31/good/Tun_OneQueue_yes.target
[Unit]
Description=Foo
[Tun]
OneQueue=yes

