# tests/auto/InvalidValue/2.31/good/Tun_User_iamauser.target
[Unit]
Description=Foo
[Tun]
User=iamauser

