# tests/auto/InvalidValue/2.31/good/Tun_OneQueue_no.target
[Unit]
Description=Foo
[Tun]
OneQueue=no

