# tests/auto/InvalidValue/2.31/good/Tun_PacketInfo_on.target
[Unit]
Description=Foo
[Tun]
PacketInfo=on

