# tests/auto/InvalidValue/2.31/good/Tun_MultiQueue_no.target
[Unit]
Description=Foo
[Tun]
MultiQueue=no

