# tests/auto/InvalidValue/2.31/good/Tun_PacketInfo_off.target
[Unit]
Description=Foo
[Tun]
PacketInfo=off

