# tests/auto/InvalidValue/2.31/good/Tun_MultiQueue_yes.target
[Unit]
Description=Foo
[Tun]
MultiQueue=yes

