# tests/auto/InvalidValue/2.31/good/Tun_PacketInfo_1.target
[Unit]
Description=Foo
[Tun]
PacketInfo=1

