# tests/auto/InvalidValue/2.31/good/Tun_OneQueue_false.target
[Unit]
Description=Foo
[Tun]
OneQueue=false

