# tests/auto/InvalidValue/2.31/good/Tun_PacketInfo_no.target
[Unit]
Description=Foo
[Tun]
PacketInfo=no

