# tests/auto/InvalidValue/2.31/good/Tunnel_ISATAP_0.target
[Unit]
Description=Foo
[Tunnel]
ISATAP=0

