# tests/auto/InvalidValue/2.31/good/Tunnel_IPv6RapidDeploymentPrefix_abc.target
[Unit]
Description=Foo
[Tunnel]
IPv6RapidDeploymentPrefix=abc

