# tests/auto/InvalidValue/2.31/good/Tunnel_Independent_0.target
[Unit]
Description=Foo
[Tunnel]
Independent=0

