# tests/auto/InvalidValue/2.31/good/Tunnel_FooOverUDP_false.target
[Unit]
Description=Foo
[Tunnel]
FooOverUDP=false

