# tests/auto/InvalidValue/2.31/good/Tunnel_SerializeTunneledPackets_on.target
[Unit]
Description=Foo
[Tunnel]
SerializeTunneledPackets=on

