# tests/auto/InvalidValue/2.31/good/Tunnel_TOS_abc.target
[Unit]
Description=Foo
[Tunnel]
TOS=abc

