# tests/auto/InvalidValue/2.31/good/Tunnel_SerializeTunneledPackets_0.target
[Unit]
Description=Foo
[Tunnel]
SerializeTunneledPackets=0

