# tests/auto/InvalidValue/2.31/good/Tunnel_ISATAP_no.target
[Unit]
Description=Foo
[Tunnel]
ISATAP=no

