# tests/auto/InvalidValue/2.31/good/Tunnel_Independent_true.target
[Unit]
Description=Foo
[Tunnel]
Independent=true

