# tests/auto/InvalidValue/2.31/good/Tunnel_DiscoverPathMTU_false.target
[Unit]
Description=Foo
[Tunnel]
DiscoverPathMTU=false

