# tests/auto/InvalidValue/2.31/good/Tunnel_FooOverUDP_1.target
[Unit]
Description=Foo
[Tunnel]
FooOverUDP=1

