# tests/auto/InvalidValue/2.31/good/Tunnel_ISATAP_1.target
[Unit]
Description=Foo
[Tunnel]
ISATAP=1

