# tests/auto/InvalidValue/2.31/good/Tunnel_DiscoverPathMTU_true.target
[Unit]
Description=Foo
[Tunnel]
DiscoverPathMTU=true

