# tests/auto/InvalidValue/2.31/good/Tunnel_FooOverUDP_yes.target
[Unit]
Description=Foo
[Tunnel]
FooOverUDP=yes

