# tests/auto/InvalidValue/2.31/good/Tunnel_SerializeTunneledPackets_no.target
[Unit]
Description=Foo
[Tunnel]
SerializeTunneledPackets=no

