# tests/auto/InvalidValue/2.31/good/Tunnel_FooOverUDP_on.target
[Unit]
Description=Foo
[Tunnel]
FooOverUDP=on

