# tests/auto/InvalidValue/2.31/good/Tunnel_Independent_no.target
[Unit]
Description=Foo
[Tunnel]
Independent=no

