# tests/auto/InvalidValue/2.31/good/Tunnel_Independent_false.target
[Unit]
Description=Foo
[Tunnel]
Independent=false

