# tests/auto/InvalidValue/2.31/good/Tunnel_FooOverUDP_true.target
[Unit]
Description=Foo
[Tunnel]
FooOverUDP=true

