# tests/auto/InvalidValue/2.31/good/Tunnel_FooOverUDP_0.target
[Unit]
Description=Foo
[Tunnel]
FooOverUDP=0

