# tests/auto/InvalidValue/2.31/good/Tunnel_Encapsulation_abc.target
[Unit]
Description=Foo
[Tunnel]
Encapsulation=abc

