# tests/auto/InvalidValue/2.31/good/Tunnel_Mode_ipip6.target
[Unit]
Description=Foo
[Tunnel]
Mode=ipip6

