# tests/auto/InvalidValue/2.31/good/Tunnel_CopyDSCP_on.target
[Unit]
Description=Foo
[Tunnel]
CopyDSCP=on

