# tests/auto/InvalidValue/2.31/good/Tunnel_Remote_255.255.255.255.target
[Unit]
Description=Foo
[Tunnel]
Remote=255.255.255.255

