# tests/auto/InvalidValue/2.31/good/Tunnel_EncapsulationLimit_255.target
[Unit]
Description=Foo
[Tunnel]
EncapsulationLimit=255

