# tests/auto/InvalidValue/2.31/good/Tunnel_Key_abc.target
[Unit]
Description=Foo
[Tunnel]
Key=abc

