# tests/auto/InvalidValue/2.31/good/Tunnel_ISATAP_true.target
[Unit]
Description=Foo
[Tunnel]
ISATAP=true

