# tests/auto/InvalidValue/2.31/good/Tunnel_FooOverUDP_no.target
[Unit]
Description=Foo
[Tunnel]
FooOverUDP=no

