# tests/auto/InvalidValue/2.31/good/Tunnel_Independent_off.target
[Unit]
Description=Foo
[Tunnel]
Independent=off

