# tests/auto/InvalidValue/2.31/good/Tunnel_CopyDSCP_yes.target
[Unit]
Description=Foo
[Tunnel]
CopyDSCP=yes

