# tests/auto/InvalidValue/2.31/good/Tunnel_Mode_any.target
[Unit]
Description=Foo
[Tunnel]
Mode=any

