# tests/auto/InvalidValue/2.31/good/Tunnel_CopyDSCP_no.target
[Unit]
Description=Foo
[Tunnel]
CopyDSCP=no

