# tests/auto/InvalidValue/2.31/good/Tunnel_ERSPANIndex_1048575.target
[Unit]
Description=Foo
[Tunnel]
ERSPANIndex=1048575

