# tests/auto/InvalidValue/2.31/good/Tunnel_AllowLocalRemote_off.target
[Unit]
Description=Foo
[Tunnel]
AllowLocalRemote=off

