# tests/auto/InvalidValue/2.31/good/Unit_ConditionACPower_true.target
[Unit]
Description=Foo
Description=abc
ConditionACPower=true

