# tests/auto/InvalidValue/2.31/good/Unit_SourcePath_abc.target
[Unit]
Description=Foo
Description=abc
SourcePath=abc

