# tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutRebootArgument_reboot-immediate.target
[Unit]
Description=Foo
Description=abc
JobTimeoutRebootArgument=reboot-immediate

