# tests/auto/InvalidValue/2.31/good/Unit_Also_abc.target
[Unit]
Description=Foo
Description=abc
Also=abc

