# tests/auto/InvalidValue/2.31/good/Unit_JobRunningTimeoutSec_1day.target
[Unit]
Description=Foo
Description=abc
JobRunningTimeoutSec=1day

