# tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_23week.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=23week

