# tests/auto/InvalidValue/2.31/good/Unit_IgnoreOnIsolate_on.target
[Unit]
Description=Foo
Description=abc
IgnoreOnIsolate=on

