# tests/auto/InvalidValue/2.31/good/Unit_JobRunningTimeoutSec_26year.target
[Unit]
Description=Foo
Description=abc
JobRunningTimeoutSec=26year

