# tests/auto/InvalidValue/2.31/good/Unit_AssertACPower_off.target
[Unit]
Description=Foo
Description=abc
AssertACPower=off

