# tests/auto/InvalidValue/2.31/good/Unit_StartLimitAction_exit-force.target
[Unit]
Description=Foo
Description=abc
StartLimitAction=exit-force

