# tests/auto/InvalidValue/2.31/good/Unit_ConditionNeedsUpdate_abc.target
[Unit]
Description=Foo
Description=abc
ConditionNeedsUpdate=abc

