# tests/auto/InvalidValue/2.31/good/Unit_FailureAction_poweroff.target
[Unit]
Description=Foo
Description=abc
FailureAction=poweroff

