# tests/auto/InvalidValue/2.31/good/Unit_AssertSecurity_abc.target
[Unit]
Description=Foo
Description=abc
AssertSecurity=abc

