# tests/auto/InvalidValue/2.31/good/Unit_RefuseManualStart_off.target
[Unit]
Description=Foo
Description=abc
RefuseManualStart=off

