# tests/auto/InvalidValue/2.31/good/Unit_RefuseManualStart_0.target
[Unit]
Description=Foo
Description=abc
RefuseManualStart=0

