# tests/auto/InvalidValue/2.31/good/Unit_IgnoreOnIsolate_false.target
[Unit]
Description=Foo
Description=abc
IgnoreOnIsolate=false

