# tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutAction_exit.target
[Unit]
Description=Foo
Description=abc
JobTimeoutAction=exit

