# tests/auto/InvalidValue/2.31/good/Unit_ConditionHost_abc.target
[Unit]
Description=Foo
Description=abc
ConditionHost=abc

