# tests/auto/InvalidValue/2.31/good/Unit_AssertPathIsMountPoint_--dev--.target
[Unit]
Description=Foo
Description=abc
AssertPathIsMountPoint=/dev/

