# tests/auto/InvalidValue/2.31/good/Unit_AssertVirtualization_false.target
[Unit]
Description=Foo
Description=abc
AssertVirtualization=false

