# tests/auto/InvalidValue/2.31/good/Unit_AssertPathExistsGlob_abc.target
[Unit]
Description=Foo
Description=abc
AssertPathExistsGlob=abc

