# tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_10M.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=10M

