# tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutRebootArgument_none.target
[Unit]
Description=Foo
Description=abc
JobTimeoutRebootArgument=none

