# tests/auto/InvalidValue/2.31/good/Unit_IgnoreOnIsolate_yes.target
[Unit]
Description=Foo
Description=abc
IgnoreOnIsolate=yes

