# tests/auto/InvalidValue/2.31/good/Unit_AssertFirstBoot_on.target
[Unit]
Description=Foo
Description=abc
AssertFirstBoot=on

