# tests/auto/InvalidValue/2.31/good/Unit_RefuseManualStop_no.target
[Unit]
Description=Foo
Description=abc
RefuseManualStop=no

