# tests/auto/InvalidValue/2.31/good/Unit_AllowIsolate_0.target
[Unit]
Description=Foo
Description=abc
AllowIsolate=0

