# tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_28%C2%B5s.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=28µs

