# tests/auto/InvalidValue/2.31/good/Unit_JoinsNamespaceOf_abc.target
[Unit]
Description=Foo
Description=abc
JoinsNamespaceOf=abc

