# tests/auto/InvalidValue/2.31/good/Unit_AllowIsolate_false.target
[Unit]
Description=Foo
Description=abc
AllowIsolate=false

