# tests/auto/InvalidValue/2.31/good/Unit_ConditionCapability_abc.target
[Unit]
Description=Foo
Description=abc
ConditionCapability=abc

