# tests/auto/InvalidValue/2.31/good/Unit_StopWhenUnneeded_true.target
[Unit]
Description=Foo
Description=abc
StopWhenUnneeded=true

