# tests/auto/InvalidValue/2.31/good/Unit_ConditionPathExistsGlob_abc.target
[Unit]
Description=Foo
Description=abc
ConditionPathExistsGlob=abc

