# tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutRebootArgument_poweroff-immediate.target
[Unit]
Description=Foo
Description=abc
JobTimeoutRebootArgument=poweroff-immediate

