# tests/auto/InvalidValue/2.31/good/Unit_AssertCPUs_0.target
[Unit]
Description=Foo
Description=abc
AssertCPUs=0

