# tests/auto/InvalidValue/2.31/good/Unit_ConditionFileNotEmpty_--dev--.target
[Unit]
Description=Foo
Description=abc
ConditionFileNotEmpty=/dev/

