# tests/auto/InvalidValue/2.31/good/Unit_ConditionVirtualization_true.target
[Unit]
Description=Foo
Description=abc
ConditionVirtualization=true

