# tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutAction_poweroff.target
[Unit]
Description=Foo
Description=abc
JobTimeoutAction=poweroff

