# tests/auto/InvalidValue/2.31/good/Unit_RefuseManualStop_on.target
[Unit]
Description=Foo
Description=abc
RefuseManualStop=on

