# tests/auto/InvalidValue/2.31/good/Unit_AssertGroup_abc.target
[Unit]
Description=Foo
Description=abc
AssertGroup=abc

