# tests/auto/InvalidValue/2.31/good/Unit_IgnoreOnIsolate_no.target
[Unit]
Description=Foo
Description=abc
IgnoreOnIsolate=no

