# tests/auto/InvalidValue/2.31/good/Unit_StartLimitIntervalSec_9999999999999.target
[Unit]
Description=Foo
Description=abc
StartLimitIntervalSec=9999999999999

