# tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_18sec.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=18sec

