# tests/auto/InvalidValue/2.31/good/Unit_AssertMemory_abc.target
[Unit]
Description=Foo
Description=abc
AssertMemory=abc

