# tests/auto/InvalidValue/2.31/good/Unit_JobRunningTimeoutSec_27y.target
[Unit]
Description=Foo
Description=abc
JobRunningTimeoutSec=27y

