# tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_24w.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=24w

