# tests/auto/InvalidValue/2.31/good/Unit_StartLimitAction_reboot.target
[Unit]
Description=Foo
Description=abc
StartLimitAction=reboot

