# tests/auto/InvalidValue/2.31/good/Unit_OnFailureJobMode_ignore-requirements.target
[Unit]
Description=Foo
Description=abc
OnFailureJobMode=ignore-requirements

