# tests/auto/InvalidValue/2.31/good/Unit_ConditionPathIsMountPoint_--dev--.target
[Unit]
Description=Foo
Description=abc
ConditionPathIsMountPoint=/dev/

