# tests/auto/InvalidValue/2.31/good/Unit_DefaultDependencies_1.target
[Unit]
Description=Foo
Description=abc
DefaultDependencies=1

