# tests/auto/InvalidValue/2.31/good/Unit_JobRunningTimeoutSec_infinity.target
[Unit]
Description=Foo
Description=abc
JobRunningTimeoutSec=infinity

