# tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutAction_reboot-immediate.target
[Unit]
Description=Foo
Description=abc
JobTimeoutAction=reboot-immediate

