# tests/auto/InvalidValue/2.31/good/Unit_IgnoreOnIsolate_off.target
[Unit]
Description=Foo
Description=abc
IgnoreOnIsolate=off

