# tests/auto/InvalidValue/2.31/good/Unit_DefaultDependencies_false.target
[Unit]
Description=Foo
Description=abc
DefaultDependencies=false

