# tests/auto/InvalidValue/2.31/good/Unit_StartLimitBurst_9999999999999.target
[Unit]
Description=Foo
Description=abc
StartLimitBurst=9999999999999

