# tests/auto/InvalidValue/2.31/good/Unit_AssertVirtualization_off.target
[Unit]
Description=Foo
Description=abc
AssertVirtualization=off

