# tests/auto/InvalidValue/2.31/good/Unit_AssertControlGroupController_abc.target
[Unit]
Description=Foo
Description=abc
AssertControlGroupController=abc

