# tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_9min.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=9min

