# tests/auto/InvalidValue/2.31/good/Unit_ConditionACPower_0.target
[Unit]
Description=Foo
Description=abc
ConditionACPower=0

