# tests/auto/InvalidValue/2.31/good/Unit_AssertUser_abc.target
[Unit]
Description=Foo
Description=abc
AssertUser=abc

