# tests/auto/InvalidValue/2.31/good/Unit_JobTimeoutSec_3hours.target
[Unit]
Description=Foo
Description=abc
JobTimeoutSec=3hours

