# tests/auto/InvalidValue/2.31/good/Unit_FailureAction_reboot-immediate.target
[Unit]
Description=Foo
Description=abc
FailureAction=reboot-immediate

