# tests/auto/InvalidValue/2.31/good/VLAN_LooseBinding_on.target
[Unit]
Description=Foo
[VLAN]
LooseBinding=on

