# tests/auto/InvalidValue/2.31/good/VLAN_MVRP_yes.target
[Unit]
Description=Foo
[VLAN]
MVRP=yes

