# tests/auto/InvalidValue/2.31/good/VLAN_LooseBinding_yes.target
[Unit]
Description=Foo
[VLAN]
LooseBinding=yes

