# tests/auto/InvalidValue/2.31/good/VLAN_ReorderHeader_on.target
[Unit]
Description=Foo
[VLAN]
ReorderHeader=on

