# tests/auto/InvalidValue/2.31/good/VLAN_ReorderHeader_true.target
[Unit]
Description=Foo
[VLAN]
ReorderHeader=true

