# tests/auto/InvalidValue/2.31/good/VLAN_LooseBinding_1.target
[Unit]
Description=Foo
[VLAN]
LooseBinding=1

