# tests/auto/InvalidValue/2.31/good/VLAN_ReorderHeader_no.target
[Unit]
Description=Foo
[VLAN]
ReorderHeader=no

