# tests/auto/InvalidValue/2.31/good/VLAN_MVRP_false.target
[Unit]
Description=Foo
[VLAN]
MVRP=false

