# tests/auto/InvalidValue/2.31/good/VLAN_LooseBinding_no.target
[Unit]
Description=Foo
[VLAN]
LooseBinding=no

