# tests/auto/InvalidValue/2.31/good/VLAN_GVRP_false.target
[Unit]
Description=Foo
[VLAN]
GVRP=false

