# tests/auto/InvalidValue/2.31/good/VLAN_GVRP_off.target
[Unit]
Description=Foo
[VLAN]
GVRP=off

