# tests/auto/InvalidValue/2.31/good/VLAN_GVRP_yes.target
[Unit]
Description=Foo
[VLAN]
GVRP=yes

