# tests/auto/InvalidValue/2.31/good/VXLAN_DestinationPort_1.target
[Unit]
Description=Foo
[VXLAN]
DestinationPort=1

