# tests/auto/InvalidValue/2.31/good/VXLAN_RemoteChecksumRx_1.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumRx=1

