# tests/auto/InvalidValue/2.31/good/VXLAN_RemoteChecksumTx_true.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumTx=true

