# tests/auto/InvalidValue/2.31/good/VXLAN_RouteShortCircuit_yes.target
[Unit]
Description=Foo
[VXLAN]
RouteShortCircuit=yes

