# tests/auto/InvalidValue/2.31/good/VXLAN_RemoteChecksumRx_true.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumRx=true

