# tests/auto/InvalidValue/2.31/good/VXLAN_MacLearning_true.target
[Unit]
Description=Foo
[VXLAN]
MacLearning=true

