# tests/auto/InvalidValue/2.31/good/VXLAN_UDPChecksum_on.target
[Unit]
Description=Foo
[VXLAN]
UDPChecksum=on

