# tests/auto/InvalidValue/2.31/good/VXLAN_RemoteChecksumTx_on.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumTx=on

