# tests/auto/InvalidValue/2.31/good/VXLAN_RemoteChecksumRx_0.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumRx=0

