# tests/auto/InvalidValue/2.31/good/VXLAN_RouteShortCircuit_true.target
[Unit]
Description=Foo
[VXLAN]
RouteShortCircuit=true

