# tests/auto/InvalidValue/2.31/good/VXLAN_RemoteChecksumTx_1.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumTx=1

