# tests/auto/InvalidValue/2.31/good/VXLAN_MacLearning_1.target
[Unit]
Description=Foo
[VXLAN]
MacLearning=1

