# tests/auto/InvalidValue/2.31/good/VXLAN_RouteShortCircuit_no.target
[Unit]
Description=Foo
[VXLAN]
RouteShortCircuit=no

