# tests/auto/InvalidValue/2.31/good/VXLAN_MaximumFDBEntries_0.target
[Unit]
Description=Foo
[VXLAN]
MaximumFDBEntries=0

