# tests/auto/InvalidValue/2.31/good/VXLAN_RemoteChecksumTx_yes.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumTx=yes

