# tests/auto/InvalidValue/2.31/good/VXLAN_L3MissNotification_on.target
[Unit]
Description=Foo
[VXLAN]
L3MissNotification=on

