# tests/auto/InvalidValue/2.31/good/VXLAN_FlowLabel_0.target
[Unit]
Description=Foo
[VXLAN]
FlowLabel=0

