# tests/auto/InvalidValue/2.31/good/VXLAN_RouteShortCircuit_false.target
[Unit]
Description=Foo
[VXLAN]
RouteShortCircuit=false

