# tests/auto/InvalidValue/2.31/good/VXLAN_RouteShortCircuit_off.target
[Unit]
Description=Foo
[VXLAN]
RouteShortCircuit=off

