# tests/auto/InvalidValue/2.31/good/VXLAN_L2MissNotification_1.target
[Unit]
Description=Foo
[VXLAN]
L2MissNotification=1

