# tests/auto/InvalidValue/2.31/good/VXLAN_PortRange_abc.target
[Unit]
Description=Foo
[VXLAN]
PortRange=abc

