# tests/auto/InvalidValue/2.31/good/VXLAN_RouteShortCircuit_1.target
[Unit]
Description=Foo
[VXLAN]
RouteShortCircuit=1

