# tests/auto/InvalidValue/2.31/good/VXLAN_ReduceARPProxy_true.target
[Unit]
Description=Foo
[VXLAN]
ReduceARPProxy=true

