# tests/auto/InvalidValue/2.31/good/VXLAN_RouteShortCircuit_on.target
[Unit]
Description=Foo
[VXLAN]
RouteShortCircuit=on

