# tests/auto/InvalidValue/2.31/good/VXLAN_UDPChecksum_no.target
[Unit]
Description=Foo
[VXLAN]
UDPChecksum=no

