# tests/auto/InvalidValue/2.31/good/VXLAN_MaximumFDBEntries_9999999999999.target
[Unit]
Description=Foo
[VXLAN]
MaximumFDBEntries=9999999999999

