# tests/auto/InvalidValue/2.31/good/VXLAN_RemoteChecksumTx_false.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumTx=false

