# tests/auto/InvalidValue/2.31/good/VXLAN_RemoteChecksumRx_on.target
[Unit]
Description=Foo
[VXLAN]
RemoteChecksumRx=on

