# tests/auto/InvalidValue/2.31/good/VXLAN_UDP6ZeroChecksumRx_0.target
[Unit]
Description=Foo
[VXLAN]
UDP6ZeroChecksumRx=0

