# tests/auto/InvalidValue/2.31/good/VXLAN_Remote_255.255.255.255--32.target
[Unit]
Description=Foo
[VXLAN]
Remote=255.255.255.255/32

