# tests/auto/InvalidValue/2.31/good/WireGuardPeer_PresharedKeyFile_--tmp.target
[Unit]
Description=Foo
[WireGuardPeer]
PresharedKeyFile=/tmp

