# tests/auto/InvalidValue/2.31/good/WireGuardPeer_PersistentKeepalive_1.target
[Unit]
Description=Foo
[WireGuardPeer]
PersistentKeepalive=1

