# tests/auto/InvalidValue/2.31/good/WireGuard_PrivateKey_abc.target
[Unit]
Description=Foo
[WireGuard]
PrivateKey=abc

