# tests/auto/InvalidValue/2.31/good/WireGuardPeer_FwMark_abc.target
[Unit]
Description=Foo
[WireGuardPeer]
FwMark=abc

