# tests/auto/InvalidValue/2.31/good/WireGuardPeer_Endpoint_abc.target
[Unit]
Description=Foo
[WireGuardPeer]
Endpoint=abc

