# tests/auto/InvalidValue/2.31/good/WireGuardPeer_PresharedKey_abc.target
[Unit]
Description=Foo
[WireGuardPeer]
PresharedKey=abc

