# tests/auto/InvalidValue/2.31/good/WireGuardPeer_PublicKey_abc.target
[Unit]
Description=Foo
[WireGuardPeer]
PublicKey=abc

