# tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_20usec.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=20usec

