# tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_2d.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=2d

