# tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_19s.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=19s

