# tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_8minute.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=8minute

