# tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_6h.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=6h

