# tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_16seconds.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=16seconds

