# tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_17second.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=17second

