# tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_11m.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=11m

