# tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_1day.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=1day

