# tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_14msec.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=14msec

