# tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_23week.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=23week

