# tests/auto/InvalidValue/2.40/good/Service_MemoryMin_9999999999999M.target
[Unit]
Description=Foo
[Service]
MemoryMin=9999999999999M

