# tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_5hr.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=5hr

