# tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_22weeks.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=22weeks

