# tests/auto/InvalidValue/2.40/good/Service_LogRateLimitIntervalSec_4hour.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=4hour

