# tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_15ms.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=15ms

