# tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_4hour.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=4hour

