# tests/auto/InvalidValue/2.40/good/Service_MemoryMin_0G.target
[Unit]
Description=Foo
[Service]
MemoryMin=0G

