# tests/auto/InvalidValue/2.40/good/Service_LogRateLimitBurst_11m.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=11m

