# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_16seconds.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=16seconds

