# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_6h.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=6h

