# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_24w.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=24w

