# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_28%C2%B5s.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=28µs

