# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_2d.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=2d

