# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_28%C2%B5s.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=28µs

