# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_7minutes.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=7minutes

