# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_14msec.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=14msec

