# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_10M.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=10M

