# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_26year.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=26year

