# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_1day.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=1day

