# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_27y.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=27y

