# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_22weeks.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=22weeks

