# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_22weeks.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=22weeks

