# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_11m.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=11m

