# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_0days.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=0days

