# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitBurst_23week.target
[Unit]
Description=Foo
[Socket]
LogRateLimitBurst=23week

