# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_15ms.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=15ms

