# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_infinity.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=infinity

