# tests/auto/InvalidValue/2.40/good/Socket_LogRateLimitIntervalSec_17second.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=17second

