# tests/auto/InvalidValue/2.40/good/Swap_MemoryMin_9999999999999B.target
[Unit]
Description=Foo
[Swap]
MemoryMin=9999999999999B

