# tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitBurst_18sec.target
[Unit]
Description=Foo
[Swap]
LogRateLimitBurst=18sec

