# tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitBurst_14msec.target
[Unit]
Description=Foo
[Swap]
LogRateLimitBurst=14msec

