# tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitBurst_4hour.target
[Unit]
Description=Foo
[Swap]
LogRateLimitBurst=4hour

