# tests/auto/InvalidValue/2.40/good/Swap_MemoryMin_0M.target
[Unit]
Description=Foo
[Swap]
MemoryMin=0M

