# tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitIntervalSec_28%C2%B5s.target
[Unit]
Description=Foo
[Swap]
LogRateLimitIntervalSec=28µs

