# tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitIntervalSec_11m.target
[Unit]
Description=Foo
[Swap]
LogRateLimitIntervalSec=11m

