# tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitIntervalSec_15ms.target
[Unit]
Description=Foo
[Swap]
LogRateLimitIntervalSec=15ms

