# tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitIntervalSec_27y.target
[Unit]
Description=Foo
[Swap]
LogRateLimitIntervalSec=27y

