# tests/auto/InvalidValue/2.40/good/Swap_MemoryMin_0K.target
[Unit]
Description=Foo
[Swap]
MemoryMin=0K

