# tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitIntervalSec_1day.target
[Unit]
Description=Foo
[Swap]
LogRateLimitIntervalSec=1day

