# tests/auto/InvalidValue/2.40/good/Swap_LogRateLimitIntervalSec_8minute.target
[Unit]
Description=Foo
[Swap]
LogRateLimitIntervalSec=8minute

