# tests/auto/InvalidValue/2.40/good/Unit_FailureActionExitStatus_255.target
[Unit]
Description=Foo
Description=abc
FailureActionExitStatus=255

