# systemdlint/tests/auto/OptionTooNew/2.44/bad/FairQueueing_Buckets_0.target
[Unit]
Description=Foo
[FairQueueing]
Buckets=123

