# systemdlint/tests/auto/OptionTooNew/2.44/bad/FairQueueing_Parent_root.target
[Unit]
Description=Foo
[FairQueueing]
Parent=123

