# systemdlint/tests/auto/OptionTooNew/2.44/bad/NetworkEmulator_DelayJitterSec_0days.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=123

