# systemdlint/tests/auto/OptionTooNew/2.44/bad/QDics_Parent_clsact.target
[Unit]
Description=Foo
[QDics]
Parent=123

