# systemdlint/tests/auto/OptionTooNew/2.44/bad/RoutingPolicyRule_User_abc.target
[Unit]
Description=Foo
[RoutingPolicyRule]
User=123

