# systemdlint/tests/auto/OptionTooNew/2.44/bad/StochasticFairnessQueueing_Handle_00.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
Handle=123

