# tests/auto/InvalidValue/2.45/good/Address_AddPrefixRoute_1.target
[Unit]
Description=Foo
[Address]
Address=abc
AddPrefixRoute=1

