# tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_3hours.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=3hours

