# tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_7minutes.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=7minutes

