# tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_5hr.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=5hr

