# tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_infinity.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=infinity

