# tests/auto/InvalidValue/2.45/good/ControlledDelay_TargetSec_28%C2%B5s.target
[Unit]
Description=Foo
[ControlledDelay]
TargetSec=28µs

