# tests/auto/InvalidValue/2.45/good/ControlledDelay_TargetSec_17second.target
[Unit]
Description=Foo
[ControlledDelay]
TargetSec=17second

