# tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_22weeks.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=22weeks

