# tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_12month.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=12month

