# tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_20usec.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=20usec

