# tests/auto/InvalidValue/2.45/good/ControlledDelay_Parent_ingress.target
[Unit]
Description=Foo
[ControlledDelay]
Parent=ingress

