# tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_10M.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=10M

