# tests/auto/InvalidValue/2.45/good/ControlledDelay_Handle_fff.target
[Unit]
Description=Foo
[ControlledDelay]
Handle=fff

