# tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_0days.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=0days

