# tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_infinity.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=infinity

