# tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_13months.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=13months

