# tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_4hour.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=4hour

