# tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_6h.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=6h

