# tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_10M.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=10M

