# tests/auto/InvalidValue/2.45/good/ControlledDelay_TargetSec_21us.target
[Unit]
Description=Foo
[ControlledDelay]
TargetSec=21us

