# tests/auto/InvalidValue/2.45/good/ControlledDelay_TargetSec_22weeks.target
[Unit]
Description=Foo
[ControlledDelay]
TargetSec=22weeks

