# tests/auto/InvalidValue/2.45/good/ControlledDelay_TargetSec_10M.target
[Unit]
Description=Foo
[ControlledDelay]
TargetSec=10M

