# tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_20usec.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=20usec

