# tests/auto/InvalidValue/2.45/good/ControlledDelay_ECN_off.target
[Unit]
Description=Foo
[ControlledDelay]
ECN=off

