# tests/auto/InvalidValue/2.45/good/ControlledDelay_TargetSec_24w.target
[Unit]
Description=Foo
[ControlledDelay]
TargetSec=24w

