# tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_24w.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=24w

