# tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_11m.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=11m

