# tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_16seconds.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=16seconds

