# tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_23week.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=23week

