# tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_5hr.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=5hr

