# tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_13months.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=13months

