# tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_25years.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=25years

