# tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_19s.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=19s

