# tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_2d.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=2d

