# tests/auto/InvalidValue/2.45/good/ControlledDelay_ECN_no.target
[Unit]
Description=Foo
[ControlledDelay]
ECN=no

