# tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_23week.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=23week

