# tests/auto/InvalidValue/2.45/good/ControlledDelay_Parent_clsact.target
[Unit]
Description=Foo
[ControlledDelay]
Parent=clsact

