# tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_16seconds.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=16seconds

