# tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_26year.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=26year

