# tests/auto/InvalidValue/2.45/good/ControlledDelay_ECN_on.target
[Unit]
Description=Foo
[ControlledDelay]
ECN=on

