# tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_24w.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=24w

