# tests/auto/InvalidValue/2.45/good/ControlledDelay_TargetSec_19s.target
[Unit]
Description=Foo
[ControlledDelay]
TargetSec=19s

