# tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_3hours.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=3hours

