# tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_8minute.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=8minute

