# tests/auto/InvalidValue/2.45/good/ControlledDelay_CEThresholdSec_18sec.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=18sec

