# tests/auto/InvalidValue/2.45/good/ControlledDelay_IntervalSec_27y.target
[Unit]
Description=Foo
[ControlledDelay]
IntervalSec=27y

