# tests/auto/InvalidValue/2.45/good/DHCPv4_SendDecline_0.target
[Unit]
Description=Foo
[DHCPv4]
SendDecline=0

