# tests/auto/InvalidValue/2.45/good/DHCPv4_SendDecline_on.target
[Unit]
Description=Foo
[DHCPv4]
SendDecline=on

