# tests/auto/InvalidValue/2.45/good/FairQueueing_Pacing_no.target
[Unit]
Description=Foo
[FairQueueing]
Pacing=no

