# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_20usec.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=20usec

