# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_10M.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=10M

