# tests/auto/InvalidValue/2.45/good/FairQueueing_Handle_fff.target
[Unit]
Description=Foo
[FairQueueing]
Handle=fff

