# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_MemoryLimit_0G.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
MemoryLimit=0G

