# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_3hours.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=3hours

