# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_Flows_9999999999999.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
Flows=9999999999999

