# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_Parent_ingress.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
Parent=ingress

