# tests/auto/InvalidValue/2.45/good/FairQueueing_MaximumRate_9999999999999G.target
[Unit]
Description=Foo
[FairQueueing]
MaximumRate=9999999999999G

