# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_11m.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=11m

