# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_18sec.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=18sec

