# tests/auto/InvalidValue/2.45/good/FairQueueing_CEThresholdSec_8minute.target
[Unit]
Description=Foo
[FairQueueing]
CEThresholdSec=8minute

