# tests/auto/InvalidValue/2.45/good/FairQueueing_Parent_ingress.target
[Unit]
Description=Foo
[FairQueueing]
Parent=ingress

