# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_1day.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=1day

