# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_19s.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=19s

