# tests/auto/InvalidValue/2.45/good/FairQueueing_CEThresholdSec_22weeks.target
[Unit]
Description=Foo
[FairQueueing]
CEThresholdSec=22weeks

