# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_19s.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=19s

