# tests/auto/InvalidValue/2.45/good/FairQueueing_CEThresholdSec_26year.target
[Unit]
Description=Foo
[FairQueueing]
CEThresholdSec=26year

