# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_ECN_1.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
ECN=1

