# tests/auto/InvalidValue/2.45/good/FairQueueing_PacketLimit_9999999999999.target
[Unit]
Description=Foo
[FairQueueing]
PacketLimit=9999999999999

