# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_13months.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=13months

