# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_Quantum_0G.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
Quantum=0G

