# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_26year.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=26year

