# tests/auto/InvalidValue/2.45/good/FairQueueing_Parent_root.target
[Unit]
Description=Foo
[FairQueueing]
Parent=root

