# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_24w.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=24w

