# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_5hr.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=5hr

