# tests/auto/InvalidValue/2.45/good/FairQueueing_MaximumRate_0M.target
[Unit]
Description=Foo
[FairQueueing]
MaximumRate=0M

