# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_15ms.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=15ms

