# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_12month.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=12month

