# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_28%C2%B5s.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=28µs

