# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_Flows_0.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
Flows=0

