# tests/auto/InvalidValue/2.45/good/FairQueueing_CEThresholdSec_24w.target
[Unit]
Description=Foo
[FairQueueing]
CEThresholdSec=24w

