# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_ECN_true.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
ECN=true

