# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_28%C2%B5s.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=28µs

