# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_17second.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=17second

