# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_13months.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=13months

