# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_26year.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=26year

