# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_25years.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=25years

