# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_Quantum_0K.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
Quantum=0K

