# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_15ms.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=15ms

