# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_Parent_clsact.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
Parent=clsact

