# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_Quantum_9999999999999G.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
Quantum=9999999999999G

