# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_21us.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=21us

