# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_PacketLimit_0.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
PacketLimit=0

