# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_16seconds.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=16seconds

