# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_8minute.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=8minute

