# tests/auto/InvalidValue/2.45/good/FairQueueing_CEThresholdSec_11m.target
[Unit]
Description=Foo
[FairQueueing]
CEThresholdSec=11m

