# tests/auto/InvalidValue/2.45/good/FairQueueing_Pacing_on.target
[Unit]
Description=Foo
[FairQueueing]
Pacing=on

