# tests/auto/InvalidValue/2.45/good/FairQueueing_CEThresholdSec_20usec.target
[Unit]
Description=Foo
[FairQueueing]
CEThresholdSec=20usec

