# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_14msec.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=14msec

