# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_infinity.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=infinity

