# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_1day.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=1day

