# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_27y.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=27y

