# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_7minutes.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=7minutes

