# tests/auto/InvalidValue/2.45/good/FairQueueing_InitialQuantum_0G.target
[Unit]
Description=Foo
[FairQueueing]
InitialQuantum=0G

