# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_12month.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=12month

