# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_16seconds.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=16seconds

