# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_17second.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=17second

