# tests/auto/InvalidValue/2.45/good/FairQueueing_OrphanMask_9999999999999.target
[Unit]
Description=Foo
[FairQueueing]
OrphanMask=9999999999999

