# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_14msec.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=14msec

