# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_TargetSec_16seconds.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=16seconds

