# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_21us.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=21us

