# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_14msec.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=14msec

