# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_MemoryLimit_9999999999999M.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
MemoryLimit=9999999999999M

