# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_ECN_off.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
ECN=off

