# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_4hour.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=4hour

