# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_ECN_no.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
ECN=no

