# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_4hour.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=4hour

