# tests/auto/InvalidValue/2.45/good/FairQueueing_Pacing_1.target
[Unit]
Description=Foo
[FairQueueing]
Pacing=1

