# tests/auto/InvalidValue/2.45/good/FairQueueing_CEThresholdSec_21us.target
[Unit]
Description=Foo
[FairQueueing]
CEThresholdSec=21us

