# tests/auto/InvalidValue/2.45/good/FairQueueing_Pacing_0.target
[Unit]
Description=Foo
[FairQueueing]
Pacing=0

