# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_22weeks.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=22weeks

