# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_20usec.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=20usec

