# tests/auto/InvalidValue/2.45/good/FairQueueing_Buckets_9999999999999.target
[Unit]
Description=Foo
[FairQueueing]
Buckets=9999999999999

