# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_24w.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=24w

