# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_infinity.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=infinity

