# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_21us.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=21us

