# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_CEThresholdSec_23week.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=23week

