# tests/auto/InvalidValue/2.45/good/FairQueueing_CEThresholdSec_13months.target
[Unit]
Description=Foo
[FairQueueing]
CEThresholdSec=13months

