# tests/auto/InvalidValue/2.45/good/FairQueueing_Quantum_9999999999999M.target
[Unit]
Description=Foo
[FairQueueing]
Quantum=9999999999999M

