# tests/auto/InvalidValue/2.45/good/FairQueueingControlledDelay_IntervalSec_8minute.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=8minute

