# tests/auto/InvalidValue/2.45/good/Link_AlternativeNamesPolicy_onboard.target
[Unit]
Description=Foo
[Link]
AlternativeNamesPolicy=onboard

