# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_19s.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=19s

