# tests/auto/InvalidValue/2.45/good/NetworkEmulator_Handle_00.target
[Unit]
Description=Foo
[NetworkEmulator]
Handle=00

