# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_3hours.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=3hours

