# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_27y.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=27y

