# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_10M.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=10M

