# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_18sec.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=18sec

