# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_1years%202months%203days%204hours%205min%2030usec.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=1years 2months 3days 4hours 5min 30usec

