# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_16seconds.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=16seconds

