# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_4hour.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=4hour

