# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_28%C2%B5s.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=28µs

