# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelaySec_9999999999999.target
[Unit]
Description=Foo
[NetworkEmulator]
DelaySec=9999999999999

