# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_0days.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=0days

