# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_17second.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=17second

