# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelaySec_0.target
[Unit]
Description=Foo
[NetworkEmulator]
DelaySec=0

