# tests/auto/InvalidValue/2.45/good/NetworkEmulator_Parent_root.target
[Unit]
Description=Foo
[NetworkEmulator]
Parent=root

