# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DuplicateRate_0%25.target
[Unit]
Description=Foo
[NetworkEmulator]
DuplicateRate=0%

