# tests/auto/InvalidValue/2.45/good/NetworkEmulator_LossRate_0%25.target
[Unit]
Description=Foo
[NetworkEmulator]
LossRate=0%

