# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_2d.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=2d

