# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_5hr.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=5hr

