# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_15ms.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=15ms

