# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_7minutes.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=7minutes

