# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_24w.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=24w

