# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_21us.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=21us

