# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_25years.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=25years

