# tests/auto/InvalidValue/2.45/good/NetworkEmulator_DelayJitterSec_13months.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=13months

