# tests/auto/InvalidValue/2.45/good/QDics_Handle_fff.target
[Unit]
Description=Foo
[QDics]
Handle=fff

