# tests/auto/InvalidValue/2.45/good/RoutingPolicyRule_SuppressPrefixLength_128.target
[Unit]
Description=Foo
[RoutingPolicyRule]
SuppressPrefixLength=128

