# tests/auto/InvalidValue/2.45/good/Service_ProtectClock_off.target
[Unit]
Description=Foo
[Service]
ProtectClock=off

