# tests/auto/InvalidValue/2.45/good/Service_ProtectClock_yes.target
[Unit]
Description=Foo
[Service]
ProtectClock=yes

