# tests/auto/InvalidValue/2.45/good/Service_ProtectClock_no.target
[Unit]
Description=Foo
[Service]
ProtectClock=no

