# tests/auto/InvalidValue/2.45/good/Socket_ProtectClock_no.target
[Unit]
Description=Foo
[Socket]
ProtectClock=no

