# tests/auto/InvalidValue/2.45/good/Socket_ProtectClock_false.target
[Unit]
Description=Foo
[Socket]
ProtectClock=false

