# tests/auto/InvalidValue/2.45/good/Socket_ProtectClock_yes.target
[Unit]
Description=Foo
[Socket]
ProtectClock=yes

