# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_18sec.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=18sec

