# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_23week.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=23week

