# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_5hr.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=5hr

