# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_4hour.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=4hour

