# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_25years.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=25years

