# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_Parent_root.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
Parent=root

