# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_14msec.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=14msec

