# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_12month.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=12month

