# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_infinity.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=infinity

