# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_7minutes.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=7minutes

