# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_13months.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=13months

