# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_2d.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=2d

