# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_Parent_clsact.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
Parent=clsact

