# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_Handle_00.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
Handle=00

