# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_24w.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=24w

