# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_27y.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=27y

