# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_20usec.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=20usec

