# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_28%C2%B5s.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=28µs

