# tests/auto/InvalidValue/2.45/good/StochasticFairnessQueueing_PerturbPeriodSec_0days.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=0days

