# tests/auto/InvalidValue/2.45/good/Swap_ProtectClock_1.target
[Unit]
Description=Foo
[Swap]
ProtectClock=1

