# tests/auto/InvalidValue/2.45/good/Swap_ProtectClock_0.target
[Unit]
Description=Foo
[Swap]
ProtectClock=0

