# tests/auto/InvalidValue/2.45/good/Swap_ProtectClock_off.target
[Unit]
Description=Foo
[Swap]
ProtectClock=off

