# tests/auto/InvalidValue/2.45/good/TokenBucketFilter_Rate_0M.target
[Unit]
Description=Foo
[TokenBucketFilter]
Rate=0M

