# tests/auto/InvalidValue/2.45/good/TokenBucketFilter_LatencySec_7minutes.target
[Unit]
Description=Foo
[TokenBucketFilter]
LatencySec=7minutes

