# tests/auto/InvalidValue/2.45/good/TokenBucketFilter_LatencySec_infinity.target
[Unit]
Description=Foo
[TokenBucketFilter]
LatencySec=infinity

