# tests/auto/InvalidValue/2.45/good/TokenBucketFilter_LatencySec_6h.target
[Unit]
Description=Foo
[TokenBucketFilter]
LatencySec=6h

