# tests/auto/InvalidValue/2.45/good/TokenBucketFilter_Rate_9999999999999K.target
[Unit]
Description=Foo
[TokenBucketFilter]
Rate=9999999999999K

