# tests/auto/InvalidValue/2.45/good/TokenBucketFilter_Parent_clsact.target
[Unit]
Description=Foo
[TokenBucketFilter]
Parent=clsact

