# tests/auto/InvalidValue/2.45/good/TokenBucketFilter_LatencySec_25years.target
[Unit]
Description=Foo
[TokenBucketFilter]
LatencySec=25years

