# tests/auto/InvalidValue/2.45/good/TokenBucketFilter_Handle_fff.target
[Unit]
Description=Foo
[TokenBucketFilter]
Handle=fff

