# tests/auto/InvalidValue/2.45/good/TokenBucketFilter_Parent_root.target
[Unit]
Description=Foo
[TokenBucketFilter]
Parent=root

