# tests/auto/InvalidValue/2.45/good/TrivialLinkEqualizer_Handle_00.target
[Unit]
Description=Foo
[TrivialLinkEqualizer]
Handle=00

