# tests/auto/InvalidValue/2.45/good/TrivialLinkEqualizer_Parent_clsact.target
[Unit]
Description=Foo
[TrivialLinkEqualizer]
Parent=clsact

