# tests/auto/InvalidValue/2.45/good/TrivialLinkEqualizer_Parent_root.target
[Unit]
Description=Foo
[TrivialLinkEqualizer]
Parent=root

