# tests/auto/OptionTooNew/2.39/good/Service_PrivateMounts_yes.target
[Unit]
Description=Foo
[Service]
PrivateMounts=123

