# systemdlint/tests/auto/OptionTooNew/2.33/bad/Service_StandardInputText_abc.target
[Unit]
Description=Foo
[Service]
StandardInputText=123

