# systemdlint/tests/auto/OptionTooNew/2.33/bad/Socket_StandardInput_null.target
[Unit]
Description=Foo
[Socket]
StandardInput=123

