# systemdlint/tests/auto/OptionTooNew/2.33/bad/Swap_StandardInputText_abc.target
[Unit]
Description=Foo
[Swap]
StandardInputText=123

