# systemdlint/tests/auto/OptionTooNew/2.33/bad/Unit_SuccessAction_none.target
[Unit]
Description=Foo
Description=abc
SuccessAction=123

