# tests/auto/OptionTooNew/2.49/good/DHCPServerStaticLease_MACAddress_01%3A23%3A45%3A67%3A89%3Aab.target
[Unit]
Description=Foo
[DHCPServerStaticLease]
MACAddress=123

