# tests/auto/OptionTooNew/2.49/good/DHCPServer_RelayTarget_1.1.1.1--0.target
[Unit]
Description=Foo
[DHCPServer]
RelayTarget=123

