# tests/auto/OptionTooNew/2.49/good/Network_ManageForeignRoutingPolicyRules_yes.target
[Unit]
Description=Foo
[Network]
ManageForeignRoutingPolicyRules=123

