# tests/auto/OptionTooNew/2.49/good/Scope_SocketBindDeny_abc.target
[Unit]
Description=Foo
[Scope]
SocketBindDeny=123

