# tests/auto/OptionTooNew/2.49/good/Socket_SocketBindDeny_abc.target
[Unit]
Description=Foo
[Socket]
SocketBindDeny=123

