# tests/auto/OptionTooNew/2.49/good/Unit_BPFProgram_abc.target
[Unit]
Description=Foo
Description=abc
BPFProgram=123

