# tests/auto/OptionTooNew/2.49/good/Unit_OnSuccess_syslog.service.target
[Unit]
Description=Foo
Description=abc
OnSuccess=123

