# tests/auto/OptionTooNew/2.49/good/Unit_ConditionOSRelease_abc.target
[Unit]
Description=Foo
Description=abc
ConditionOSRelease=123

