# tests/auto/InvalidValue/2.35/bad/Mount_IPAddressDeny_0.0.0..target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
IPAddressDeny=0.0.0.

