# tests/auto/InvalidValue/2.35/bad/Mount_IPAddressDeny_256.255.255.255.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
IPAddressDeny=256.255.255.255

