# tests/auto/InvalidValue/2.35/bad/Mount_IPAddressAllow_0--abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
IPAddressAllow=0/abc

