# tests/auto/InvalidValue/2.35/bad/Mount_IPAddressAllow_-1.1.1.1.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
IPAddressAllow=-1.1.1.1

