# tests/auto/InvalidValue/2.35/bad/Service_StateDirectoryMode_10000000000000.target
[Unit]
Description=Foo
[Service]
StateDirectoryMode=10000000000000

