# tests/auto/InvalidValue/2.35/bad/Service_IPAddressAllow_abc.target
[Unit]
Description=Foo
[Service]
IPAddressAllow=abc

