# tests/auto/InvalidValue/2.35/bad/Service_IPAddressAllow_0.0.0..target
[Unit]
Description=Foo
[Service]
IPAddressAllow=0.0.0.

