# tests/auto/InvalidValue/2.35/bad/Service_StateDirectoryMode_-1.target
[Unit]
Description=Foo
[Service]
StateDirectoryMode=-1

