# tests/auto/InvalidValue/2.35/bad/Service_IPAddressDeny_1.1.1.1---1.target
[Unit]
Description=Foo
[Service]
IPAddressDeny=1.1.1.1/-1

