# tests/auto/InvalidValue/2.35/bad/Service_CacheDirectoryMode_10000000000000.target
[Unit]
Description=Foo
[Service]
CacheDirectoryMode=10000000000000

