# tests/auto/InvalidValue/2.35/bad/Service_IPAddressAllow_-1.1.1.1.target
[Unit]
Description=Foo
[Service]
IPAddressAllow=-1.1.1.1

