# tests/auto/InvalidValue/2.35/bad/Service_IPAddressDeny_255.256.255.255.target
[Unit]
Description=Foo
[Service]
IPAddressDeny=255.256.255.255

