# tests/auto/InvalidValue/2.35/bad/Service_ConfigurationDirectoryMode_10000000000000.target
[Unit]
Description=Foo
[Service]
ConfigurationDirectoryMode=10000000000000

