# tests/auto/InvalidValue/2.35/bad/Socket_IPAddressDeny_-1.1.1.1.target
[Unit]
Description=Foo
[Socket]
IPAddressDeny=-1.1.1.1

