# tests/auto/InvalidValue/2.35/bad/Socket_IPAddressDeny_abc.target
[Unit]
Description=Foo
[Socket]
IPAddressDeny=abc

