# tests/auto/InvalidValue/2.35/bad/Socket_IPAddressAllow_0.0.0..target
[Unit]
Description=Foo
[Socket]
IPAddressAllow=0.0.0.

