# tests/auto/InvalidValue/2.35/bad/Socket_IPAccounting_foo.target
[Unit]
Description=Foo
[Socket]
IPAccounting=foo

