# tests/auto/InvalidValue/2.35/bad/Socket_IPAddressDeny_0.0.0..target
[Unit]
Description=Foo
[Socket]
IPAddressDeny=0.0.0.

