# tests/auto/InvalidValue/2.35/bad/Socket_IPAddressDeny_255.256.255.255.target
[Unit]
Description=Foo
[Socket]
IPAddressDeny=255.256.255.255

