# tests/auto/InvalidValue/2.35/bad/Socket_IPAddressAllow_0--abc.target
[Unit]
Description=Foo
[Socket]
IPAddressAllow=0/abc

