# tests/auto/InvalidValue/2.35/bad/Socket_ConfigurationDirectoryMode_-1.target
[Unit]
Description=Foo
[Socket]
ConfigurationDirectoryMode=-1

