# tests/auto/InvalidValue/2.35/bad/Socket_RuntimeDirectoryMode_10000000000000.target
[Unit]
Description=Foo
[Socket]
RuntimeDirectoryMode=10000000000000

