# tests/auto/InvalidValue/2.35/bad/Socket_IPAddressAllow_255.256.255.255.target
[Unit]
Description=Foo
[Socket]
IPAddressAllow=255.256.255.255

