# tests/auto/InvalidValue/2.35/bad/Socket_LockPersonality_baz.target
[Unit]
Description=Foo
[Socket]
LockPersonality=baz

