# tests/auto/InvalidValue/2.35/bad/Socket_IPAddressAllow_abc.target
[Unit]
Description=Foo
[Socket]
IPAddressAllow=abc

