# tests/auto/InvalidValue/2.35/bad/Socket_IPAddressDeny_256.255.255.255.target
[Unit]
Description=Foo
[Socket]
IPAddressDeny=256.255.255.255

