# tests/auto/InvalidValue/2.35/bad/Swap_IPAccounting_foo.target
[Unit]
Description=Foo
[Swap]
IPAccounting=foo

