# tests/auto/InvalidValue/2.35/bad/Swap_LockPersonality_foo.target
[Unit]
Description=Foo
[Swap]
LockPersonality=foo

