# tests/auto/InvalidValue/2.35/bad/Swap_LockPersonality_bar.target
[Unit]
Description=Foo
[Swap]
LockPersonality=bar

