# tests/auto/InvalidValue/2.35/bad/Swap_CacheDirectoryMode_10000000000000.target
[Unit]
Description=Foo
[Swap]
CacheDirectoryMode=10000000000000

