# tests/auto/InvalidValue/2.35/bad/Swap_IPAccounting_baz.target
[Unit]
Description=Foo
[Swap]
IPAccounting=baz

