# tests/auto/InvalidValue/2.35/bad/Swap_LockPersonality_baz.target
[Unit]
Description=Foo
[Swap]
LockPersonality=baz

