# tests/auto/InvalidValue/2.35/bad/Swap_IPAddressDeny_255.256.255.255.target
[Unit]
Description=Foo
[Swap]
IPAddressDeny=255.256.255.255

