# tests/auto/InvalidValue/2.35/bad/Swap_StateDirectoryMode_10000000000000.target
[Unit]
Description=Foo
[Swap]
StateDirectoryMode=10000000000000

