# tests/auto/InvalidValue/2.35/bad/Swap_IPAddressDeny_abc.target
[Unit]
Description=Foo
[Swap]
IPAddressDeny=abc

