# tests/auto/InvalidValue/2.35/bad/Swap_IPAddressAllow_abc.target
[Unit]
Description=Foo
[Swap]
IPAddressAllow=abc

