# tests/auto/InvalidValue/2.35/bad/Swap_IPAddressAllow_0--abc.target
[Unit]
Description=Foo
[Swap]
IPAddressAllow=0/abc

