# tests/auto/InvalidValue/2.35/bad/Swap_IPAddressAllow_-1.1.1.1.target
[Unit]
Description=Foo
[Swap]
IPAddressAllow=-1.1.1.1

