# tests/auto/InvalidValue/2.35/bad/Swap_RuntimeDirectoryMode_10000000000000.target
[Unit]
Description=Foo
[Swap]
RuntimeDirectoryMode=10000000000000

