# tests/auto/OptionDeprecated/2.51/bad/DHCPv6PrefixDelegation_ManageTemporaryAddress_abc.target
[Unit]
Description=Foo
[DHCPv6PrefixDelegation]
ManageTemporaryAddress=123

