# tests/auto/InvalidValue/2.27/good/User_CrashChVT_abc.target
[Unit]
Description=Foo
[User]
CrashChVT=abc

