# systemdlint/tests/manual/ConflictingOptions/2.31/bad/Unit_BeforeAfter.target
[Unit]
Description=Before and After
After=foo.service
Before=foo.service