# systemdlint/tests/manual/ConflictingOptions/2.31/bad/Unit_RequiresRequiredBy.target
[Unit]
Description=Before and After
Requires=foo.service
RequiredBy=foo.service