# tests/auto/OptionTooNew/2.47/good/Mount_ManagedOOMMemoryPressure_abc.target
[Unit]
Description=Foo
[Mount]
What=abc
Where=/tmp
ManagedOOMMemoryPressure=123

