# tests/auto/OptionTooNew/2.47/good/Service_ManagedOOMMemoryPressureLimitPercent_abc.target
[Unit]
Description=Foo
[Service]
ManagedOOMMemoryPressureLimitPercent=123

