# tests/auto/OptionTooNew/2.47/good/Service_SystemCallLog_abc.target
[Unit]
Description=Foo
[Service]
SystemCallLog=123

