# tests/auto/OptionTooNew/2.47/good/Service_ManagedOOMSwap_abc.target
[Unit]
Description=Foo
[Service]
ManagedOOMSwap=123

