# tests/auto/OptionTooNew/2.47/good/Service_ManagedOOMMemoryPressure_abc.target
[Unit]
Description=Foo
[Service]
ManagedOOMMemoryPressure=123

