# tests/auto/OptionTooNew/2.47/good/Socket_NUMAMask_all.target
[Unit]
Description=Foo
[Socket]
NUMAMask=123

