# tests/auto/OptionTooNew/2.47/good/Swap_ManagedOOMSwap_abc.target
[Unit]
Description=Foo
[Swap]
ManagedOOMSwap=123

