# systemdlint/tests/auto/SettingRestricted/2.31/good/DHCP_ClientIdentifier_macAnonymize%3Dtrue.target
[Unit]
Description=Foo
[DHCP]
ClientIdentifier=mac

