# systemdlint/tests/auto/SettingRestricted/2.31/good/DHCP_VendorClassIdentifier_abcAnonymize%3D1.target
[Unit]
Description=Foo
[DHCP]
VendorClassIdentifier=abc

