# systemdlint/tests/auto/SettingRestricted/2.31/good/Network_ActiveSlave_yesMode%3Dactive-backup.target
[Unit]
Description=Foo
[Network]
ActiveSlave=yes

