# systemdlint/tests/auto/SettingRestricted/2.31/good/Path_MakeDirectory_yesPathExists%3DNone.target
[Unit]
Description=Foo
[Path]
MakeDirectory=yes

