# systemdlint/tests/auto/MandatoryOptionMissing/2.31/good/Unit_Mandatory_Unit.target
[Unit]
Description=Foo
[Unit]
Description=abc
