# tests/auto/OptionDeprecated/2.27/good/User_CrashChVT_abc.target
[Unit]
Description=Foo
[User]
CrashChVT=123

