# systemdlint/tests/auto/OptionDeprecated/2.12/good/Service_ControlGroupPersistent_abc.target
[Unit]
Description=Foo
[Service]
ControlGroupPersistent=123

