# systemdlint/tests/auto/OptionDeprecated/2.12/good/Service_ControlGroup_abc.target
[Unit]
Description=Foo
[Service]
ControlGroup=123

