# systemdlint/tests/auto/OptionDeprecated/2.12/good/Socket_ControlGroupPersistent_abc.target
[Unit]
Description=Foo
[Socket]
ControlGroupPersistent=123

