# systemdlint/tests/auto/OptionDeprecated/2.12/good/Swap_ControlGroup_abc.target
[Unit]
Description=Foo
[Swap]
ControlGroup=123

