# systemdlint/tests/auto/InvalidValue/2.39/good/Manager_NoNewPrivileges_true.target
[Unit]
Description=Foo
[Manager]
NoNewPrivileges=true

