# systemdlint/tests/auto/InvalidValue/2.39/good/Manager_NoNewPrivileges_on.target
[Unit]
Description=Foo
[Manager]
NoNewPrivileges=on

