# systemdlint/tests/auto/InvalidValue/2.39/good/Manager_NoNewPrivileges_no.target
[Unit]
Description=Foo
[Manager]
NoNewPrivileges=no

