# systemdlint/tests/auto/InvalidValue/2.39/good/Manager_NoNewPrivileges_off.target
[Unit]
Description=Foo
[Manager]
NoNewPrivileges=off

