# systemdlint/tests/auto/InvalidValue/2.39/good/Manager_NoNewPrivileges_yes.target
[Unit]
Description=Foo
[Manager]
NoNewPrivileges=yes

