# systemdlint/tests/auto/InvalidValue/2.39/good/Manager_NoNewPrivileges_false.target
[Unit]
Description=Foo
[Manager]
NoNewPrivileges=false

