# systemdlint/tests/auto/InvalidValue/2.39/good/Service_PrivateMounts_false.target
[Unit]
Description=Foo
[Service]
PrivateMounts=false

