# systemdlint/tests/auto/InvalidValue/2.39/good/Service_PrivateMounts_on.target
[Unit]
Description=Foo
[Service]
PrivateMounts=on

