# systemdlint/tests/auto/InvalidValue/2.39/good/Service_NoNewPrivileges_true.target
[Unit]
Description=Foo
[Service]
NoNewPrivileges=true

