# systemdlint/tests/auto/InvalidValue/2.39/good/Service_PrivateMounts_no.target
[Unit]
Description=Foo
[Service]
PrivateMounts=no

