# systemdlint/tests/auto/InvalidValue/2.39/good/Service_PrivateMounts_off.target
[Unit]
Description=Foo
[Service]
PrivateMounts=off

