# systemdlint/tests/auto/InvalidValue/2.39/good/Socket_NoNewPrivileges_yes.target
[Unit]
Description=Foo
[Socket]
NoNewPrivileges=yes

