# systemdlint/tests/auto/InvalidValue/2.39/good/Socket_PrivateMounts_yes.target
[Unit]
Description=Foo
[Socket]
PrivateMounts=yes

