# systemdlint/tests/auto/InvalidValue/2.39/good/Socket_NoNewPrivileges_true.target
[Unit]
Description=Foo
[Socket]
NoNewPrivileges=true

