# systemdlint/tests/auto/InvalidValue/2.39/good/Socket_PrivateMounts_false.target
[Unit]
Description=Foo
[Socket]
PrivateMounts=false

