# systemdlint/tests/auto/InvalidValue/2.39/good/Socket_PrivateMounts_off.target
[Unit]
Description=Foo
[Socket]
PrivateMounts=off

