# systemdlint/tests/auto/InvalidValue/2.39/good/Socket_NoNewPrivileges_on.target
[Unit]
Description=Foo
[Socket]
NoNewPrivileges=on

