# systemdlint/tests/auto/InvalidValue/2.39/good/Socket_NoNewPrivileges_off.target
[Unit]
Description=Foo
[Socket]
NoNewPrivileges=off

