# systemdlint/tests/auto/InvalidValue/2.39/good/Socket_NoNewPrivileges_1.target
[Unit]
Description=Foo
[Socket]
NoNewPrivileges=1

