# systemdlint/tests/auto/InvalidValue/2.39/good/Socket_PrivateMounts_1.target
[Unit]
Description=Foo
[Socket]
PrivateMounts=1

