# systemdlint/tests/auto/InvalidValue/2.39/good/Socket_NoNewPrivileges_no.target
[Unit]
Description=Foo
[Socket]
NoNewPrivileges=no

