# systemdlint/tests/auto/InvalidValue/2.39/good/Swap_PrivateMounts_true.target
[Unit]
Description=Foo
[Swap]
PrivateMounts=true

