# systemdlint/tests/auto/InvalidValue/2.39/good/Swap_PrivateMounts_yes.target
[Unit]
Description=Foo
[Swap]
PrivateMounts=yes

