# systemdlint/tests/auto/InvalidValue/2.39/good/Swap_PrivateMounts_on.target
[Unit]
Description=Foo
[Swap]
PrivateMounts=on

