# systemdlint/tests/auto/InvalidValue/2.39/good/Swap_NoNewPrivileges_on.target
[Unit]
Description=Foo
[Swap]
NoNewPrivileges=on

