# systemdlint/tests/auto/InvalidValue/2.39/good/Swap_NoNewPrivileges_no.target
[Unit]
Description=Foo
[Swap]
NoNewPrivileges=no

