# systemdlint/tests/auto/InvalidValue/2.39/good/Swap_PrivateMounts_1.target
[Unit]
Description=Foo
[Swap]
PrivateMounts=1

