# systemdlint/tests/auto/InvalidValue/2.39/good/Swap_PrivateMounts_off.target
[Unit]
Description=Foo
[Swap]
PrivateMounts=off

