# tests/auto/InvalidValue/2.40/bad/Service_LogRateLimitBurst_True.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=True

