# tests/auto/InvalidValue/2.40/bad/Service_MemoryMin_10000000000000M.target
[Unit]
Description=Foo
[Service]
MemoryMin=10000000000000M

