# tests/auto/InvalidValue/2.40/bad/Service_LogRateLimitBurst_0.5d.target
[Unit]
Description=Foo
[Service]
LogRateLimitBurst=0.5d

