# tests/auto/InvalidValue/2.40/bad/Service_LogRateLimitIntervalSec_2potatoes.target
[Unit]
Description=Foo
[Service]
LogRateLimitIntervalSec=2potatoes

