# tests/auto/InvalidValue/2.40/bad/Socket_LogRateLimitIntervalSec_2potatoes.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=2potatoes

