# tests/auto/InvalidValue/2.40/bad/Socket_MemoryMin_10000000000000T.target
[Unit]
Description=Foo
[Socket]
MemoryMin=10000000000000T

