# tests/auto/InvalidValue/2.40/bad/Socket_LogRateLimitIntervalSec_0.5d.target
[Unit]
Description=Foo
[Socket]
LogRateLimitIntervalSec=0.5d

