# tests/auto/InvalidValue/2.40/bad/Socket_MemoryMin_-1T.target
[Unit]
Description=Foo
[Socket]
MemoryMin=-1T

