# tests/auto/InvalidValue/2.40/bad/Swap_MemoryMin_10000000000000G.target
[Unit]
Description=Foo
[Swap]
MemoryMin=10000000000000G

