# tests/auto/InvalidValue/2.40/bad/Swap_MemoryMin_-1B.target
[Unit]
Description=Foo
[Swap]
MemoryMin=-1B

