# tests/auto/InvalidValue/2.40/bad/Swap_LogRateLimitBurst_2potatoes.target
[Unit]
Description=Foo
[Swap]
LogRateLimitBurst=2potatoes

