# tests/auto/InvalidValue/2.40/bad/Swap_MemoryMin_-1K.target
[Unit]
Description=Foo
[Swap]
MemoryMin=-1K

