# tests/auto/InvalidValue/2.40/bad/Swap_LogRateLimitBurst_True.target
[Unit]
Description=Foo
[Swap]
LogRateLimitBurst=True

