# tests/auto/OptionTooNew/2.20/good/Automount_TimeoutIdleSec_0days.target
[Unit]
Description=Foo
[Automount]
TimeoutIdleSec=123

