# tests/auto/InvalidValue/2.45/bad/ControlledDelay_CEThresholdSec_2potatoes.target
[Unit]
Description=Foo
[ControlledDelay]
CEThresholdSec=2potatoes

