# tests/auto/InvalidValue/2.45/bad/ControlledDelay_TargetSec_1p.target
[Unit]
Description=Foo
[ControlledDelay]
TargetSec=1p

