# tests/auto/InvalidValue/2.45/bad/ControlledDelay_TargetSec_3decades.target
[Unit]
Description=Foo
[ControlledDelay]
TargetSec=3decades

