# tests/auto/InvalidValue/2.45/bad/ControlledDelay_PacketLimit_4294967295.target
[Unit]
Description=Foo
[ControlledDelay]
PacketLimit=4294967295

