# tests/auto/InvalidValue/2.45/bad/FairQueueingControlledDelay_IntervalSec_0.5d.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=0.5d

