# tests/auto/InvalidValue/2.45/bad/FairQueueing_Quantum_10000000000000G.target
[Unit]
Description=Foo
[FairQueueing]
Quantum=10000000000000G

