# tests/auto/InvalidValue/2.45/bad/FairQueueingControlledDelay_ECN_bar.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
ECN=bar

