# tests/auto/InvalidValue/2.45/bad/FairQueueingControlledDelay_TargetSec_True.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=True

