# tests/auto/InvalidValue/2.45/bad/FairQueueing_Pacing_bar.target
[Unit]
Description=Foo
[FairQueueing]
Pacing=bar

