# tests/auto/InvalidValue/2.45/bad/FairQueueingControlledDelay_PacketLimit_-1.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
PacketLimit=-1

