# tests/auto/InvalidValue/2.45/bad/FairQueueingControlledDelay_MemoryLimit_10000000000000M.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
MemoryLimit=10000000000000M

