# tests/auto/InvalidValue/2.45/bad/FairQueueing_PacketLimit_-1.target
[Unit]
Description=Foo
[FairQueueing]
PacketLimit=-1

