# tests/auto/InvalidValue/2.45/bad/FairQueueingControlledDelay_CEThresholdSec_0.5d.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=0.5d

