# tests/auto/InvalidValue/2.45/bad/FairQueueingControlledDelay_ECN_foo.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
ECN=foo

