# tests/auto/InvalidValue/2.45/bad/FairQueueingControlledDelay_ECN_baz.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
ECN=baz

