# tests/auto/InvalidValue/2.45/bad/FairQueueingControlledDelay_MemoryLimit_-1K.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
MemoryLimit=-1K

