# tests/auto/InvalidValue/2.45/bad/FairQueueing_CEThresholdSec_True.target
[Unit]
Description=Foo
[FairQueueing]
CEThresholdSec=True

