# tests/auto/InvalidValue/2.45/bad/FairQueueingControlledDelay_CEThresholdSec_True.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
CEThresholdSec=True

