# tests/auto/InvalidValue/2.45/bad/FairQueueing_InitialQuantum_-1G.target
[Unit]
Description=Foo
[FairQueueing]
InitialQuantum=-1G

