# tests/auto/InvalidValue/2.45/bad/FairQueueingControlledDelay_IntervalSec_True.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=True

