# tests/auto/InvalidValue/2.45/bad/FairQueueingControlledDelay_TargetSec_1p.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
TargetSec=1p

