# tests/auto/InvalidValue/2.45/bad/FairQueueing_Buckets_-1.target
[Unit]
Description=Foo
[FairQueueing]
Buckets=-1

