# tests/auto/InvalidValue/2.45/bad/FairQueueingControlledDelay_IntervalSec_3decades.target
[Unit]
Description=Foo
[FairQueueingControlledDelay]
IntervalSec=3decades

