# tests/auto/InvalidValue/2.45/bad/NetworkEmulator_DelayJitterSec_3decades.target
[Unit]
Description=Foo
[NetworkEmulator]
DelayJitterSec=3decades

