# tests/auto/InvalidValue/2.45/bad/NetworkEmulator_PacketLimit_-1.target
[Unit]
Description=Foo
[NetworkEmulator]
PacketLimit=-1

