# tests/auto/InvalidValue/2.45/bad/NetworkEmulator_LossRate_101%25.target
[Unit]
Description=Foo
[NetworkEmulator]
LossRate=101%

