# tests/auto/InvalidValue/2.45/bad/NetworkEmulator_DelaySec_10000000000000.target
[Unit]
Description=Foo
[NetworkEmulator]
DelaySec=10000000000000

