# tests/auto/InvalidValue/2.45/bad/NetworkEmulator_PacketLimit_4294967295.target
[Unit]
Description=Foo
[NetworkEmulator]
PacketLimit=4294967295

