# tests/auto/InvalidValue/2.45/bad/Service_ProtectClock_baz.target
[Unit]
Description=Foo
[Service]
ProtectClock=baz

