# tests/auto/InvalidValue/2.45/bad/Service_ProtectClock_foo.target
[Unit]
Description=Foo
[Service]
ProtectClock=foo

