# tests/auto/InvalidValue/2.45/bad/Socket_ProtectClock_bar.target
[Unit]
Description=Foo
[Socket]
ProtectClock=bar

