# tests/auto/InvalidValue/2.45/bad/Socket_ProtectClock_baz.target
[Unit]
Description=Foo
[Socket]
ProtectClock=baz

