# tests/auto/InvalidValue/2.45/bad/StochasticFairnessQueueing_PerturbPeriodSec_2potatoes.target
[Unit]
Description=Foo
[StochasticFairnessQueueing]
PerturbPeriodSec=2potatoes

