# tests/auto/InvalidValue/2.45/bad/Swap_ProtectClock_bar.target
[Unit]
Description=Foo
[Swap]
ProtectClock=bar

