# tests/auto/InvalidValue/2.45/bad/TokenBucketFilter_Rate_-1K.target
[Unit]
Description=Foo
[TokenBucketFilter]
Rate=-1K

