# tests/auto/InvalidValue/2.45/bad/TokenBucketFilter_LatencySec_True.target
[Unit]
Description=Foo
[TokenBucketFilter]
LatencySec=True

