# tests/auto/InvalidValue/2.45/bad/TokenBucketFilter_Rate_10000000000000K.target
[Unit]
Description=Foo
[TokenBucketFilter]
Rate=10000000000000K

