# tests/auto/InvalidValue/2.45/bad/TokenBucketFilter_MPUBytes_-1M.target
[Unit]
Description=Foo
[TokenBucketFilter]
MPUBytes=-1M

