# tests/auto/InvalidValue/2.45/bad/TokenBucketFilter_LatencySec_0.5d.target
[Unit]
Description=Foo
[TokenBucketFilter]
LatencySec=0.5d

