# tests/auto/InvalidValue/2.45/bad/TokenBucketFilter_LimitSize_-1K.target
[Unit]
Description=Foo
[TokenBucketFilter]
LimitSize=-1K

