# tests/auto/InvalidValue/2.45/bad/TokenBucketFilter_Burst_-1G.target
[Unit]
Description=Foo
[TokenBucketFilter]
Burst=-1G

