# tests/auto/InvalidValue/2.45/bad/TrivialLinkEqualizer_Id_-1.target
[Unit]
Description=Foo
[TrivialLinkEqualizer]
Id=-1

