# tests/auto/OptionTooNew/2.27/good/Swap_TasksMax_1.target
[Unit]
Description=Foo
[Swap]
TasksMax=123

