# tests/auto/OptionTooNew/2.11/good/Manager_DefaultCPUAccounting_yes.target
[Unit]
Description=Foo
[Manager]
DefaultCPUAccounting=123

