# tests/auto/OptionTooNew/2.11/good/Service_RestrictAddressFamilies_abc.target
[Unit]
Description=Foo
[Service]
RestrictAddressFamilies=123

