# systemdlint/tests/auto/InvalidValue/2.48/bad/Link_TransmitQueues_0.target
[Unit]
Description=Foo
[Link]
TransmitQueues=0

