# systemdlint/tests/auto/InvalidValue/2.48/bad/Link_TransmitQueueLength_-1.target
[Unit]
Description=Foo
[Link]
TransmitQueueLength=-1

