# systemdlint/tests/auto/InvalidValue/2.48/bad/MACVLAN_BroadcastMulticastQueueLength_-1.target
[Unit]
Description=Foo
[MACVLAN]
BroadcastMulticastQueueLength=-1

