# systemdlint/tests/auto/InvalidValue/2.48/bad/MACVLAN_BroadcastMulticastQueueLength_4294967295.target
[Unit]
Description=Foo
[MACVLAN]
BroadcastMulticastQueueLength=4294967295

