# systemdlint/tests/auto/InvalidValue/2.48/bad/Manager_ManagerEnvironment_True.target
[Unit]
Description=Foo
[Manager]
ManagerEnvironment=True

