# systemdlint/tests/auto/InvalidValue/2.48/bad/Manager_ManagerEnvironment_a%3Eb.target
[Unit]
Description=Foo
[Manager]
ManagerEnvironment=a>b

