# systemdlint/tests/auto/InvalidValue/2.48/bad/OOM_DefaultMemoryPressureDurationSec_10000000000000.target
[Unit]
Description=Foo
[OOM]
DefaultMemoryPressureDurationSec=10000000000000

