# systemdlint/tests/auto/InvalidValue/2.48/bad/OOM_DefaultMemoryPressureDurationSec_-1.target
[Unit]
Description=Foo
[OOM]
DefaultMemoryPressureDurationSec=-1

