# systemdlint/tests/auto/InvalidValue/2.48/bad/Time_ConnectionRetrySec_0.target
[Unit]
Description=Foo
[Time]
ConnectionRetrySec=0

