# tests/auto/OptionTooNew/2.31/bad/Service_ProtectKernelTunables_yes.target
[Unit]
Description=Foo
[Service]
ProtectKernelTunables=123

