# tests/auto/OptionTooNew/2.31/bad/Service_DynamicUser_yes.target
[Unit]
Description=Foo
[Service]
DynamicUser=123

