# tests/auto/OptionTooNew/2.31/bad/Socket_CPUAccounting_yes.target
[Unit]
Description=Foo
[Socket]
CPUAccounting=123

