# tests/auto/OptionTooNew/2.31/bad/Socket_DynamicUser_yes.target
[Unit]
Description=Foo
[Socket]
DynamicUser=123

