# tests/manual/ConflictingOptions/2.31/bad/Unit_PropagatesReloadToReloadPropagatedFrom.target
[Unit]
Description=Before and After
PropagatesReloadTo=foo.service
ReloadPropagatedFrom=foo.service