# tests/auto/InvalidValue/2.47/bad/BridgeMDB_MulticastGroupAddress_0.0.0..target
[Unit]
Description=Foo
[BridgeMDB]
MulticastGroupAddress=0.0.0.

