# tests/auto/InvalidValue/2.47/bad/BridgeMDB_MulticastGroupAddress_abc.target
[Unit]
Description=Foo
[BridgeMDB]
MulticastGroupAddress=abc

