# tests/auto/InvalidValue/2.47/bad/FlowQueuePIE_PacketLimit_4294967295.target
[Unit]
Description=Foo
[FlowQueuePIE]
PacketLimit=4294967295

