# tests/auto/InvalidValue/2.47/bad/FlowQueuePIE_Handle_0x10000.target
[Unit]
Description=Foo
[FlowQueuePIE]
Handle=0x10000

