# tests/auto/InvalidValue/2.47/bad/IPv6SendRA_Managed_bar.target
[Unit]
Description=Foo
[IPv6SendRA]
Managed=bar

