# tests/auto/InvalidValue/2.47/bad/IPv6SendRA_OtherInformation_bar.target
[Unit]
Description=Foo
[IPv6SendRA]
OtherInformation=bar

