# tests/auto/InvalidValue/2.47/bad/IPv6SendRA_DNSLifetimeSec_True.target
[Unit]
Description=Foo
[IPv6SendRA]
DNSLifetimeSec=True

