# tests/auto/InvalidValue/2.47/bad/Socket_FlushPending_baz.target
[Unit]
Description=Foo
[Socket]
FlushPending=baz

