# tests/auto/InvalidValue/2.47/bad/VXLAN_Independent_baz.target
[Unit]
Description=Foo
[VXLAN]
Independent=baz

