# systemdlint/tests/auto/InvalidValue/2.46/good/Bridge_VLANProtocol_802.1ad.target
[Unit]
Description=Foo
[Bridge]
VLANProtocol=802.1ad

