# systemdlint/tests/auto/InvalidValue/2.46/good/CAKE_OverheadBytes_256.target
[Unit]
Description=Foo
[CAKE]
OverheadBytes=256

