# systemdlint/tests/auto/InvalidValue/2.46/good/CAKE_Parent_clsact.target
[Unit]
Description=Foo
[CAKE]
Parent=clsact

