# systemdlint/tests/auto/InvalidValue/2.46/good/CAKE_Parent_root.target
[Unit]
Description=Foo
[CAKE]
Parent=root

