# systemdlint/tests/auto/InvalidValue/2.46/good/CAKE_Bandwidth_0G.target
[Unit]
Description=Foo
[CAKE]
Bandwidth=0G

