# systemdlint/tests/auto/InvalidValue/2.46/good/CAKE_Bandwidth_9999999999999K.target
[Unit]
Description=Foo
[CAKE]
Bandwidth=9999999999999K

