# systemdlint/tests/auto/InvalidValue/2.46/good/CAN_FDMode_true.target
[Unit]
Description=Foo
[CAN]
FDMode=true

