# systemdlint/tests/auto/InvalidValue/2.46/good/CAN_DataBitRate_1M.target
[Unit]
Description=Foo
[CAN]
DataBitRate=1M

