# systemdlint/tests/auto/InvalidValue/2.46/good/CAN_Termination_1.target
[Unit]
Description=Foo
[CAN]
Termination=1

