# systemdlint/tests/auto/InvalidValue/2.46/good/CAN_ListenOnly_on.target
[Unit]
Description=Foo
[CAN]
ListenOnly=on

