# systemdlint/tests/auto/InvalidValue/2.46/good/CAN_FDNonISO_no.target
[Unit]
Description=Foo
[CAN]
FDNonISO=no

