# systemdlint/tests/auto/InvalidValue/2.46/good/CAN_FDNonISO_1.target
[Unit]
Description=Foo
[CAN]
FDNonISO=1

