# systemdlint/tests/auto/InvalidValue/2.46/good/CAN_Termination_no.target
[Unit]
Description=Foo
[CAN]
Termination=no

