# systemdlint/tests/auto/InvalidValue/2.46/good/CAN_FDMode_off.target
[Unit]
Description=Foo
[CAN]
FDMode=off

