# systemdlint/tests/auto/InvalidValue/2.46/good/CAN_Termination_0.target
[Unit]
Description=Foo
[CAN]
Termination=0

