# systemdlint/tests/auto/InvalidValue/2.46/good/CAN_ListenOnly_off.target
[Unit]
Description=Foo
[CAN]
ListenOnly=off

