# systemdlint/tests/auto/InvalidValue/2.46/good/CAN_Termination_false.target
[Unit]
Description=Foo
[CAN]
Termination=false

