# systemdlint/tests/auto/InvalidValue/2.46/good/CAN_FDNonISO_yes.target
[Unit]
Description=Foo
[CAN]
FDNonISO=yes

