# systemdlint/tests/auto/InvalidValue/2.46/good/CAN_ListenOnly_yes.target
[Unit]
Description=Foo
[CAN]
ListenOnly=yes

