# systemdlint/tests/auto/InvalidValue/2.46/good/DeficitRoundRobinScheduler_Parent_clsact.target
[Unit]
Description=Foo
[DeficitRoundRobinScheduler]
Parent=clsact

