# systemdlint/tests/auto/InvalidValue/2.46/good/DeficitRoundRobinSchedulerClass_Parent_root.target
[Unit]
Description=Foo
[DeficitRoundRobinSchedulerClass]
Parent=root

