# systemdlint/tests/auto/InvalidValue/2.46/good/DeficitRoundRobinScheduler_Handle_0xffff.target
[Unit]
Description=Foo
[DeficitRoundRobinScheduler]
Handle=0xffff

