# systemdlint/tests/auto/InvalidValue/2.46/good/DeficitRoundRobinScheduler_Parent_root.target
[Unit]
Description=Foo
[DeficitRoundRobinScheduler]
Parent=root

