# systemdlint/tests/auto/InvalidValue/2.46/good/DeficitRoundRobinScheduler_Parent_ingress.target
[Unit]
Description=Foo
[DeficitRoundRobinScheduler]
Parent=ingress

