# systemdlint/tests/auto/InvalidValue/2.46/good/DeficitRoundRobinSchedulerClass_QuantumBytes_0K.target
[Unit]
Description=Foo
[DeficitRoundRobinSchedulerClass]
QuantumBytes=0K

