# systemdlint/tests/auto/InvalidValue/2.46/good/DeficitRoundRobinScheduler_Parent_abc.target
[Unit]
Description=Foo
[DeficitRoundRobinScheduler]
Parent=abc

