# systemdlint/tests/auto/InvalidValue/2.46/good/DHCPServer_EmitLPR_off.target
[Unit]
Description=Foo
[DHCPServer]
EmitLPR=off

