# systemdlint/tests/auto/InvalidValue/2.46/good/DHCPv6PrefixDelegation_Token_abc.target
[Unit]
Description=Foo
[DHCPv6PrefixDelegation]
Token=abc

