# systemdlint/tests/auto/InvalidValue/2.46/good/DHCPServer_EmitLPR_true.target
[Unit]
Description=Foo
[DHCPServer]
EmitLPR=true

