# systemdlint/tests/auto/InvalidValue/2.46/good/DHCPv4_FallbackLeaseLifetimeSec_infinity.target
[Unit]
Description=Foo
[DHCPv4]
FallbackLeaseLifetimeSec=infinity

