# systemdlint/tests/auto/InvalidValue/2.46/good/DHCPServer_EmitPOP3_0.target
[Unit]
Description=Foo
[DHCPServer]
EmitPOP3=0

