# systemdlint/tests/auto/InvalidValue/2.46/good/DHCPServer_EmitSMTP_on.target
[Unit]
Description=Foo
[DHCPServer]
EmitSMTP=on

