# systemdlint/tests/auto/InvalidValue/2.46/good/DHCPServer_EmitSMTP_false.target
[Unit]
Description=Foo
[DHCPServer]
EmitSMTP=false

