# systemdlint/tests/auto/InvalidValue/2.46/good/DHCPServer_EmitLPR_no.target
[Unit]
Description=Foo
[DHCPServer]
EmitLPR=no

