# systemdlint/tests/auto/InvalidValue/2.46/good/DHCPv6_RequestOptions_abc.target
[Unit]
Description=Foo
[DHCPv6]
RequestOptions=abc

