# systemdlint/tests/auto/InvalidValue/2.46/good/DHCPServer_EmitSMTP_off.target
[Unit]
Description=Foo
[DHCPServer]
EmitSMTP=off

