# systemdlint/tests/auto/InvalidValue/2.46/good/DHCPv6_VendorClass_abc.target
[Unit]
Description=Foo
[DHCPv6]
VendorClass=abc

