# systemdlint/tests/auto/InvalidValue/2.46/good/DHCPServer_EmitSMTP_yes.target
[Unit]
Description=Foo
[DHCPServer]
EmitSMTP=yes

