# systemdlint/tests/auto/InvalidValue/2.46/good/EnhancedTransmissionSelection_PriorityMap_abc.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
PriorityMap=abc

