# systemdlint/tests/auto/InvalidValue/2.46/good/EnhancedTransmissionSelection_Parent_abc.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
Parent=abc

