# systemdlint/tests/auto/InvalidValue/2.46/good/EnhancedTransmissionSelection_Parent_clsact.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
Parent=clsact

