# systemdlint/tests/auto/InvalidValue/2.46/good/EnhancedTransmissionSelection_QuantumBytes_9999999999999G.target
[Unit]
Description=Foo
[EnhancedTransmissionSelection]
QuantumBytes=9999999999999G

