# systemdlint/tests/auto/InvalidValue/2.46/good/GenericRandomEarlyDetection_GenericRIO_off.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
GenericRIO=off

