# systemdlint/tests/auto/InvalidValue/2.46/good/GenericRandomEarlyDetection_VirtualQueues_16.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
VirtualQueues=16

