# systemdlint/tests/auto/InvalidValue/2.46/good/GenericRandomEarlyDetection_DefaultVirtualQueue_16.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
DefaultVirtualQueue=16

