# systemdlint/tests/auto/InvalidValue/2.46/good/GenericRandomEarlyDetection_DefaultVirtualQueue_1.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
DefaultVirtualQueue=1

