# systemdlint/tests/auto/InvalidValue/2.46/good/GenericRandomEarlyDetection_Parent_ingress.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
Parent=ingress

