# systemdlint/tests/auto/InvalidValue/2.46/good/GenericRandomEarlyDetection_Parent_clsact.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
Parent=clsact

