# systemdlint/tests/auto/InvalidValue/2.46/good/GenericRandomEarlyDetection_GenericRIO_0.target
[Unit]
Description=Foo
[GenericRandomEarlyDetection]
GenericRIO=0

