# systemdlint/tests/auto/InvalidValue/2.46/good/HeavyHitterFilter_Parent_ingress.target
[Unit]
Description=Foo
[HeavyHitterFilter]
Parent=ingress

