# systemdlint/tests/auto/InvalidValue/2.46/good/HeavyHitterFilter_PacketLimit_0.target
[Unit]
Description=Foo
[HeavyHitterFilter]
PacketLimit=0

