# systemdlint/tests/auto/InvalidValue/2.46/good/HeavyHitterFilter_Parent_root.target
[Unit]
Description=Foo
[HeavyHitterFilter]
Parent=root

