# systemdlint/tests/auto/InvalidValue/2.46/good/HierarchyTokenBucket_Parent_ingress.target
[Unit]
Description=Foo
[HierarchyTokenBucket]
Parent=ingress

