# systemdlint/tests/auto/InvalidValue/2.46/good/HierarchyTokenBucketClass_MTUBytes_0G.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
MTUBytes=0G

