# systemdlint/tests/auto/InvalidValue/2.46/good/HierarchyTokenBucketClass_OverheadBytes_9999999999999G.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
OverheadBytes=9999999999999G

