# systemdlint/tests/auto/InvalidValue/2.46/good/HierarchyTokenBucketClass_BufferBytes_0K.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
BufferBytes=0K

