# systemdlint/tests/auto/InvalidValue/2.46/good/HierarchyTokenBucket_DefaultClass_0x9184e729fff.target
[Unit]
Description=Foo
[HierarchyTokenBucket]
DefaultClass=0x9184e729fff

