# systemdlint/tests/auto/InvalidValue/2.46/good/HierarchyTokenBucketClass_CeilBufferBytes_0G.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
CeilBufferBytes=0G

