# systemdlint/tests/auto/InvalidValue/2.46/good/HierarchyTokenBucket_RateToQuantum_9999999999999.target
[Unit]
Description=Foo
[HierarchyTokenBucket]
RateToQuantum=9999999999999

