# systemdlint/tests/auto/InvalidValue/2.46/good/HierarchyTokenBucketClass_BufferBytes_0G.target
[Unit]
Description=Foo
[HierarchyTokenBucketClass]
BufferBytes=0G

