# systemdlint/tests/auto/InvalidValue/2.46/good/IPv6Prefix_Assign_no.target
[Unit]
Description=Foo
[IPv6Prefix]
Assign=no

