# systemdlint/tests/auto/InvalidValue/2.46/good/IPv6AcceptRA_DenyList_abc.target
[Unit]
Description=Foo
[IPv6AcceptRA]
DenyList=abc

