# systemdlint/tests/auto/InvalidValue/2.46/good/IPv6Prefix_Assign_off.target
[Unit]
Description=Foo
[IPv6Prefix]
Assign=off

