# systemdlint/tests/auto/InvalidValue/2.46/good/IPv6Prefix_Assign_false.target
[Unit]
Description=Foo
[IPv6Prefix]
Assign=false

