# systemdlint/tests/auto/InvalidValue/2.46/good/Link_TxFlowControl_false.target
[Unit]
Description=Foo
[Link]
TxFlowControl=false

