# systemdlint/tests/auto/InvalidValue/2.46/good/Link_TxFlowControl_true.target
[Unit]
Description=Foo
[Link]
TxFlowControl=true

