# systemdlint/tests/auto/InvalidValue/2.46/good/Link_AutoNegotiationFlowControl_no.target
[Unit]
Description=Foo
[Link]
AutoNegotiationFlowControl=no

