# systemdlint/tests/auto/InvalidValue/2.46/good/Link_RxFlowControl_1.target
[Unit]
Description=Foo
[Link]
RxFlowControl=1

