# systemdlint/tests/auto/InvalidValue/2.46/good/Link_RxFlowControl_yes.target
[Unit]
Description=Foo
[Link]
RxFlowControl=yes

