# systemdlint/tests/auto/InvalidValue/2.46/good/Link_TxFlowControl_yes.target
[Unit]
Description=Foo
[Link]
TxFlowControl=yes

