# systemdlint/tests/auto/InvalidValue/2.46/good/Link_RxFlowControl_on.target
[Unit]
Description=Foo
[Link]
RxFlowControl=on

