# systemdlint/tests/auto/InvalidValue/2.46/good/Link_RxFlowControl_true.target
[Unit]
Description=Foo
[Link]
RxFlowControl=true

