# systemdlint/tests/auto/InvalidValue/2.46/good/Link_AutoNegotiationFlowControl_1.target
[Unit]
Description=Foo
[Link]
MACAddress=abc
AutoNegotiationFlowControl=1

