# systemdlint/tests/auto/InvalidValue/2.46/good/Link_AutoNegotiationFlowControl_yes.target
[Unit]
Description=Foo
[Link]
AutoNegotiationFlowControl=yes

