# systemdlint/tests/auto/InvalidValue/2.46/good/Link_Group_4294967294.target
[Unit]
Description=Foo
[Link]
Group=4294967294

