# systemdlint/tests/auto/InvalidValue/2.46/good/Network_IPv6LinkLocalAddressGenerationMode_eui64.target
[Unit]
Description=Foo
[Network]
IPv6LinkLocalAddressGenerationMode=eui64

