# systemdlint/tests/auto/InvalidValue/2.46/good/Network_IPv6LinkLocalAddressGenerationMode_stable-privacy.target
[Unit]
Description=Foo
[Network]
IPv6LinkLocalAddressGenerationMode=stable-privacy

