# systemdlint/tests/auto/InvalidValue/2.46/good/Network_IPv4AcceptLocal_false.target
[Unit]
Description=Foo
[Network]
IPv4AcceptLocal=false

