# systemdlint/tests/auto/InvalidValue/2.46/good/PFIFOHeadDrop_PacketLimit_4294967294.target
[Unit]
Description=Foo
[PFIFOHeadDrop]
PacketLimit=4294967294

