# systemdlint/tests/auto/InvalidValue/2.46/good/PFIFO_Handle_0xffff.target
[Unit]
Description=Foo
[PFIFO]
Handle=0xffff

