# systemdlint/tests/auto/InvalidValue/2.46/good/PFIFOHeadDrop_Parent_abc.target
[Unit]
Description=Foo
[PFIFOHeadDrop]
Parent=abc

