# systemdlint/tests/auto/InvalidValue/2.46/good/PFIFOFast_Parent_root.target
[Unit]
Description=Foo
[PFIFOFast]
Parent=root

