# systemdlint/tests/auto/InvalidValue/2.46/good/PFIFOFast_Parent_clsact.target
[Unit]
Description=Foo
[PFIFOFast]
Parent=clsact

