# systemdlint/tests/auto/InvalidValue/2.46/good/PFIFOHeadDrop_Parent_root.target
[Unit]
Description=Foo
[PFIFOHeadDrop]
Parent=root

