# systemdlint/tests/auto/InvalidValue/2.46/good/PFIFOFast_Handle_0xffff.target
[Unit]
Description=Foo
[PFIFOFast]
Handle=0xffff

