# systemdlint/tests/auto/InvalidValue/2.46/good/PFIFOHeadDrop_Parent_ingress.target
[Unit]
Description=Foo
[PFIFOHeadDrop]
Parent=ingress

