# systemdlint/tests/auto/InvalidValue/2.46/good/PFIFOHeadDrop_Parent_clsact.target
[Unit]
Description=Foo
[PFIFOHeadDrop]
Parent=clsact

