# systemdlint/tests/auto/InvalidValue/2.46/good/PFIFOFast_Parent_abc.target
[Unit]
Description=Foo
[PFIFOFast]
Parent=abc

