# systemdlint/tests/auto/InvalidValue/2.46/good/PIE_PacketLimit_0.target
[Unit]
Description=Foo
[PIE]
PacketLimit=0

