# systemdlint/tests/auto/InvalidValue/2.46/good/PIE_PacketLimit_4294967294.target
[Unit]
Description=Foo
[PIE]
PacketLimit=4294967294

