# systemdlint/tests/auto/InvalidValue/2.46/good/PIE_Parent_root.target
[Unit]
Description=Foo
[PIE]
Parent=root

