# systemdlint/tests/auto/InvalidValue/2.46/good/QuickFairQueueing_Parent_clsact.target
[Unit]
Description=Foo
[QuickFairQueueing]
Parent=clsact

