# systemdlint/tests/auto/InvalidValue/2.46/good/QuickFairQueueing_Handle_0xffff.target
[Unit]
Description=Foo
[QuickFairQueueing]
Handle=0xffff

