# systemdlint/tests/auto/InvalidValue/2.46/good/QuickFairQueueingClass_MaxPacketBytes_9999999999999G.target
[Unit]
Description=Foo
[QuickFairQueueingClass]
MaxPacketBytes=9999999999999G

