# systemdlint/tests/auto/InvalidValue/2.46/good/QuickFairQueueing_Parent_root.target
[Unit]
Description=Foo
[QuickFairQueueing]
Parent=root

