# systemdlint/tests/auto/InvalidValue/2.46/good/QuickFairQueueing_Parent_ingress.target
[Unit]
Description=Foo
[QuickFairQueueing]
Parent=ingress

