# systemdlint/tests/auto/InvalidValue/2.46/good/QuickFairQueueing_Parent_abc.target
[Unit]
Description=Foo
[QuickFairQueueing]
Parent=abc

