# systemdlint/tests/auto/InvalidValue/2.46/good/QuickFairQueueingClass_ClassId_abc.target
[Unit]
Description=Foo
[QuickFairQueueingClass]
ClassId=abc

