# systemdlint/tests/auto/InvalidValue/2.46/good/QuickFairQueueingClass_MaxPacketBytes_0G.target
[Unit]
Description=Foo
[QuickFairQueueingClass]
MaxPacketBytes=0G

