# systemdlint/tests/auto/InvalidValue/2.46/good/QuickFairQueueingClass_Parent_abc.target
[Unit]
Description=Foo
[QuickFairQueueingClass]
Parent=abc

