# systemdlint/tests/auto/InvalidValue/2.46/good/Service_CoredumpFilter_all.target
[Unit]
Description=Foo
[Service]
CoredumpFilter=all

