# systemdlint/tests/auto/InvalidValue/2.46/good/Service_TimeoutStartFailureMode_terminate.target
[Unit]
Description=Foo
[Service]
TimeoutStartFailureMode=terminate

