# systemdlint/tests/auto/InvalidValue/2.46/good/Service_TimeoutStopFailureMode_terminate.target
[Unit]
Description=Foo
[Service]
TimeoutStopFailureMode=terminate

