# systemdlint/tests/auto/InvalidValue/2.46/good/Service_CoredumpFilter_private-anonymous.target
[Unit]
Description=Foo
[Service]
CoredumpFilter=private-anonymous

