# systemdlint/tests/auto/InvalidValue/2.46/good/Service_RootHash_abc.target
[Unit]
Description=Foo
[Service]
RootHash=abc

