# systemdlint/tests/auto/InvalidValue/2.46/good/Service_CoredumpFilter_default.target
[Unit]
Description=Foo
[Service]
CoredumpFilter=default

