# systemdlint/tests/auto/InvalidValue/2.46/good/Service_CoredumpFilter_shared-huge.target
[Unit]
Description=Foo
[Service]
CoredumpFilter=shared-huge

