# systemdlint/tests/auto/InvalidValue/2.46/good/Service_TimeoutStopFailureMode_abort.target
[Unit]
Description=Foo
[Service]
TimeoutStopFailureMode=abort

