# systemdlint/tests/auto/InvalidValue/2.46/good/Service_TimeoutStartFailureMode_abort.target
[Unit]
Description=Foo
[Service]
TimeoutStartFailureMode=abort

