# systemdlint/tests/auto/InvalidValue/2.46/good/Service_CoredumpFilter_private-dax.target
[Unit]
Description=Foo
[Service]
CoredumpFilter=private-dax

