# systemdlint/tests/auto/InvalidValue/2.46/good/Service_CoredumpFilter_shared-anonymous.target
[Unit]
Description=Foo
[Service]
CoredumpFilter=shared-anonymous

