# systemdlint/tests/auto/InvalidValue/2.46/good/Socket_CoredumpFilter_private-file-backed.target
[Unit]
Description=Foo
[Socket]
CoredumpFilter=private-file-backed

