# systemdlint/tests/auto/InvalidValue/2.46/good/Socket_CoredumpFilter_private-huge.target
[Unit]
Description=Foo
[Socket]
CoredumpFilter=private-huge

