# systemdlint/tests/auto/InvalidValue/2.46/good/Socket_CoredumpFilter_16.target
[Unit]
Description=Foo
[Socket]
CoredumpFilter=16

