# systemdlint/tests/auto/InvalidValue/2.46/good/Socket_CoredumpFilter_private-dax.target
[Unit]
Description=Foo
[Socket]
CoredumpFilter=private-dax

