# systemdlint/tests/auto/InvalidValue/2.46/good/Socket_RootHash_abc.target
[Unit]
Description=Foo
[Socket]
RootHash=abc

