# systemdlint/tests/auto/InvalidValue/2.46/good/Socket_CoredumpFilter_0.target
[Unit]
Description=Foo
[Socket]
CoredumpFilter=0

