# systemdlint/tests/auto/InvalidValue/2.46/good/SR-IOV_LinkState_yes.target
[Unit]
Description=Foo
[SR-IOV]
LinkState=yes

