# systemdlint/tests/auto/InvalidValue/2.46/good/SR-IOV_Trust_off.target
[Unit]
Description=Foo
[SR-IOV]
Trust=off

