# systemdlint/tests/auto/InvalidValue/2.46/good/StochasticFairBlue_Parent_abc.target
[Unit]
Description=Foo
[StochasticFairBlue]
Parent=abc

