# systemdlint/tests/auto/InvalidValue/2.46/good/Swap_CoredumpFilter_shared-anonymous.target
[Unit]
Description=Foo
[Swap]
CoredumpFilter=shared-anonymous

