# systemdlint/tests/auto/InvalidValue/2.46/good/Swap_CoredumpFilter_default.target
[Unit]
Description=Foo
[Swap]
CoredumpFilter=default

