# systemdlint/tests/auto/InvalidValue/2.46/good/Swap_CoredumpFilter_0.target
[Unit]
Description=Foo
[Swap]
CoredumpFilter=0

