# systemdlint/tests/auto/InvalidValue/2.46/good/Swap_RootHash_abc.target
[Unit]
Description=Foo
[Swap]
RootHash=abc

