# systemdlint/tests/auto/InvalidValue/2.46/good/Unit_ConditionEnvironment_a%3Db.target
[Unit]
Description=Foo
Description=abc
ConditionEnvironment=a=b

