# systemdlint/tests/auto/InvalidValue/2.46/good/Unit_AssertPathIsEncrypted_--tmp.target
[Unit]
Description=Foo
Description=abc
AssertPathIsEncrypted=/tmp

