# tests/auto/InvalidValue/2.32/bad/Service_ProtectKernelModules_baz.target
[Unit]
Description=Foo
[Service]
ProtectKernelModules=baz

