# tests/auto/InvalidValue/2.32/bad/Service_CPUWeight_0.target
[Unit]
Description=Foo
[Service]
CPUWeight=0

