# tests/auto/InvalidValue/2.32/bad/Service_DynamicUser_foo.target
[Unit]
Description=Foo
[Service]
DynamicUser=foo

