# tests/auto/InvalidValue/2.32/bad/Service_ProtectControlGroups_foo.target
[Unit]
Description=Foo
[Service]
ProtectControlGroups=foo

