# tests/auto/InvalidValue/2.32/bad/Service_PrivateUsers_bar.target
[Unit]
Description=Foo
[Service]
PrivateUsers=bar

