# tests/auto/InvalidValue/2.32/bad/Service_ProtectControlGroups_baz.target
[Unit]
Description=Foo
[Service]
ProtectControlGroups=baz

