# tests/auto/InvalidValue/2.32/bad/Service_ProtectKernelModules_foo.target
[Unit]
Description=Foo
[Service]
ProtectKernelModules=foo

