# tests/auto/InvalidValue/2.32/bad/Service_PrivateUsers_baz.target
[Unit]
Description=Foo
[Service]
PrivateUsers=baz

