# tests/auto/InvalidValue/2.32/bad/Service_DynamicUser_bar.target
[Unit]
Description=Foo
[Service]
DynamicUser=bar

