# tests/auto/InvalidValue/2.32/bad/Socket_DynamicUser_bar.target
[Unit]
Description=Foo
[Socket]
DynamicUser=bar

