# tests/auto/InvalidValue/2.32/bad/Socket_MaxConnections_10000000000000.target
[Unit]
Description=Foo
[Socket]
MaxConnections=10000000000000

