# tests/auto/InvalidValue/2.32/bad/Socket_CPUAccounting_baz.target
[Unit]
Description=Foo
[Socket]
CPUAccounting=baz

