# tests/auto/InvalidValue/2.32/bad/Socket_PrivateUsers_baz.target
[Unit]
Description=Foo
[Socket]
PrivateUsers=baz

