# tests/auto/InvalidValue/2.32/bad/Socket_ProtectControlGroups_baz.target
[Unit]
Description=Foo
[Socket]
ProtectControlGroups=baz

