# tests/auto/InvalidValue/2.32/bad/Socket_CPUWeight_10001.target
[Unit]
Description=Foo
[Socket]
CPUWeight=10001

