# tests/auto/InvalidValue/2.32/bad/Socket_MaxConnectionsPerSource_-1.target
[Unit]
Description=Foo
[Socket]
MaxConnectionsPerSource=-1

