# tests/auto/InvalidValue/2.32/bad/Socket_PrivateUsers_bar.target
[Unit]
Description=Foo
[Socket]
PrivateUsers=bar

