# tests/auto/InvalidValue/2.32/bad/Socket_CPUWeight_0.target
[Unit]
Description=Foo
[Socket]
CPUWeight=0

