# tests/auto/InvalidValue/2.32/bad/Socket_DynamicUser_baz.target
[Unit]
Description=Foo
[Socket]
DynamicUser=baz

