# tests/auto/InvalidValue/2.32/bad/Socket_ProtectKernelModules_foo.target
[Unit]
Description=Foo
[Socket]
ProtectKernelModules=foo

