# tests/auto/InvalidValue/2.32/bad/Socket_CPUAccounting_bar.target
[Unit]
Description=Foo
[Socket]
CPUAccounting=bar

