# tests/auto/InvalidValue/2.32/bad/Swap_ProtectControlGroups_bar.target
[Unit]
Description=Foo
[Swap]
ProtectControlGroups=bar

