# tests/auto/InvalidValue/2.32/bad/Swap_PrivateUsers_bar.target
[Unit]
Description=Foo
[Swap]
PrivateUsers=bar

