# tests/auto/InvalidValue/2.32/bad/Swap_ProtectKernelTunables_foo.target
[Unit]
Description=Foo
[Swap]
ProtectKernelTunables=foo

