# tests/auto/InvalidValue/2.32/bad/Swap_CPUAccounting_baz.target
[Unit]
Description=Foo
[Swap]
CPUAccounting=baz

