# tests/auto/InvalidValue/2.32/bad/Swap_ProtectKernelModules_bar.target
[Unit]
Description=Foo
[Swap]
ProtectKernelModules=bar

