# tests/auto/InvalidValue/2.32/bad/Swap_ProtectKernelModules_baz.target
[Unit]
Description=Foo
[Swap]
ProtectKernelModules=baz

