# systemdlint/tests/auto/InvalidValue/2.32/bad/Service_ProtectKernelTunables_baz.target
[Unit]
Description=Foo
[Service]
ProtectKernelTunables=baz

