# systemdlint/tests/auto/InvalidValue/2.32/bad/Service_ProtectKernelModules_bar.target
[Unit]
Description=Foo
[Service]
ProtectKernelModules=bar

