# systemdlint/tests/auto/InvalidValue/2.32/bad/Service_PrivateUsers_foo.target
[Unit]
Description=Foo
[Service]
PrivateUsers=foo

