# systemdlint/tests/auto/InvalidValue/2.32/bad/Service_CPUWeight_10001.target
[Unit]
Description=Foo
[Service]
CPUWeight=10001

