# systemdlint/tests/auto/InvalidValue/2.32/bad/Socket_MaxConnections_-1.target
[Unit]
Description=Foo
[Socket]
MaxConnections=-1

