# systemdlint/tests/auto/InvalidValue/2.32/bad/Socket_CPUAccounting_foo.target
[Unit]
Description=Foo
[Socket]
CPUAccounting=foo

