# systemdlint/tests/auto/InvalidValue/2.32/bad/Socket_ProtectControlGroups_bar.target
[Unit]
Description=Foo
[Socket]
ProtectControlGroups=bar

