# systemdlint/tests/auto/InvalidValue/2.32/bad/Socket_MaxConnectionsPerSource_10000000000000.target
[Unit]
Description=Foo
[Socket]
MaxConnectionsPerSource=10000000000000

