# systemdlint/tests/auto/InvalidValue/2.32/bad/Socket_ProtectKernelTunables_baz.target
[Unit]
Description=Foo
[Socket]
ProtectKernelTunables=baz

