# systemdlint/tests/auto/InvalidValue/2.32/bad/Socket_ProtectKernelModules_baz.target
[Unit]
Description=Foo
[Socket]
ProtectKernelModules=baz

