# systemdlint/tests/auto/InvalidValue/2.32/bad/Socket_MemorySwapMax_-1M.target
[Unit]
Description=Foo
[Socket]
MemorySwapMax=-1M

