# systemdlint/tests/auto/InvalidValue/2.32/bad/Socket_ProtectControlGroups_foo.target
[Unit]
Description=Foo
[Socket]
ProtectControlGroups=foo

