# systemdlint/tests/auto/InvalidValue/2.32/bad/Socket_ProtectKernelTunables_foo.target
[Unit]
Description=Foo
[Socket]
ProtectKernelTunables=foo

