# systemdlint/tests/auto/InvalidValue/2.32/bad/Socket_MemorySwapMax_-1G.target
[Unit]
Description=Foo
[Socket]
MemorySwapMax=-1G

