# systemdlint/tests/auto/InvalidValue/2.32/bad/Socket_ProtectKernelTunables_bar.target
[Unit]
Description=Foo
[Socket]
ProtectKernelTunables=bar

