# systemdlint/tests/auto/InvalidValue/2.32/bad/Socket_DynamicUser_foo.target
[Unit]
Description=Foo
[Socket]
DynamicUser=foo

