# systemdlint/tests/auto/InvalidValue/2.32/bad/Swap_ProtectControlGroups_foo.target
[Unit]
Description=Foo
[Swap]
ProtectControlGroups=foo

