# systemdlint/tests/auto/InvalidValue/2.32/bad/Swap_DynamicUser_bar.target
[Unit]
Description=Foo
[Swap]
DynamicUser=bar

