# systemdlint/tests/auto/InvalidValue/2.32/bad/Swap_DynamicUser_baz.target
[Unit]
Description=Foo
[Swap]
DynamicUser=baz

