# systemdlint/tests/auto/InvalidValue/2.32/bad/Swap_ProtectControlGroups_baz.target
[Unit]
Description=Foo
[Swap]
ProtectControlGroups=baz

