# systemdlint/tests/auto/InvalidValue/2.32/bad/Swap_ProtectKernelTunables_bar.target
[Unit]
Description=Foo
[Swap]
ProtectKernelTunables=bar

