# systemdlint/tests/auto/InvalidValue/2.32/bad/Swap_DynamicUser_foo.target
[Unit]
Description=Foo
[Swap]
DynamicUser=foo

