# systemdlint/tests/auto/InvalidValue/2.32/bad/Swap_ProtectKernelModules_foo.target
[Unit]
Description=Foo
[Swap]
ProtectKernelModules=foo

