# systemdlint/tests/auto/InvalidValue/2.32/bad/Swap_CPUAccounting_bar.target
[Unit]
Description=Foo
[Swap]
CPUAccounting=bar

