# systemdlint/tests/auto/InvalidValue/2.32/bad/Swap_PrivateUsers_baz.target
[Unit]
Description=Foo
[Swap]
PrivateUsers=baz

