# systemdlint/tests/auto/InvalidValue/2.32/bad/Swap_ProtectKernelTunables_baz.target
[Unit]
Description=Foo
[Swap]
ProtectKernelTunables=baz

