# systemdlint/tests/auto/InvalidValue/2.32/bad/Swap_CPUAccounting_foo.target
[Unit]
Description=Foo
[Swap]
CPUAccounting=foo

