# tests/auto/OptionTooNew/2.13/good/Unit_FailureAction_none.target
[Unit]
Description=Foo
Description=abc
FailureAction=123

