# tests/auto/OptionTooNew/2.13/good/Unit_RebootArgument_abc.target
[Unit]
Description=Foo
Description=abc
RebootArgument=123

