# tests/auto/InvalidValue/2.39/bad/Manager_NoNewPrivileges_foo.target
[Unit]
Description=Foo
[Manager]
NoNewPrivileges=foo

