# tests/auto/InvalidValue/2.39/bad/Manager_NoNewPrivileges_bar.target
[Unit]
Description=Foo
[Manager]
NoNewPrivileges=bar

