# tests/auto/InvalidValue/2.39/bad/Service_PrivateMounts_baz.target
[Unit]
Description=Foo
[Service]
PrivateMounts=baz

