# tests/auto/InvalidValue/2.39/bad/Service_PrivateMounts_bar.target
[Unit]
Description=Foo
[Service]
PrivateMounts=bar

