# tests/auto/InvalidValue/2.39/bad/Socket_PrivateMounts_foo.target
[Unit]
Description=Foo
[Socket]
PrivateMounts=foo

