# tests/auto/InvalidValue/2.39/bad/Socket_PrivateMounts_bar.target
[Unit]
Description=Foo
[Socket]
PrivateMounts=bar

