# tests/auto/InvalidValue/2.39/bad/Swap_PrivateMounts_baz.target
[Unit]
Description=Foo
[Swap]
PrivateMounts=baz

