# tests/auto/InvalidValue/2.39/bad/Swap_PrivateMounts_bar.target
[Unit]
Description=Foo
[Swap]
PrivateMounts=bar

