# tests/auto/InvalidValue/2.39/bad/Swap_NoNewPrivileges_foo.target
[Unit]
Description=Foo
[Swap]
NoNewPrivileges=foo

