# tests/auto/InvalidValue/2.39/bad/Swap_NoNewPrivileges_bar.target
[Unit]
Description=Foo
[Swap]
NoNewPrivileges=bar

