# systemdlint/tests/auto/InvalidValue/2.44/good/DHCPServer_EmitSIP_on.target
[Unit]
Description=Foo
[DHCPServer]
EmitSIP=on

