# systemdlint/tests/auto/InvalidValue/2.44/good/DHCPServer_EmitSIP_no.target
[Unit]
Description=Foo
[DHCPServer]
EmitSIP=no

