# systemdlint/tests/auto/InvalidValue/2.44/good/DHCPServer_SIP_abc.target
[Unit]
Description=Foo
[DHCPServer]
SIP=abc

