# systemdlint/tests/auto/InvalidValue/2.44/good/DHCPv6_PrefixDelegationHint_255.255.255.255.target
[Unit]
Description=Foo
[DHCPv6]
PrefixDelegationHint=255.255.255.255

