# systemdlint/tests/auto/InvalidValue/2.44/good/DHCPServer_EmitSIP_off.target
[Unit]
Description=Foo
[DHCPServer]
EmitSIP=off

