# systemdlint/tests/auto/InvalidValue/2.44/good/Link_TxBufferSize_0.target
[Unit]
Description=Foo
[Link]
TxBufferSize=0

