# systemdlint/tests/auto/InvalidValue/2.44/good/Link_TxBufferSize_9999999999999.target
[Unit]
Description=Foo
[Link]
TxBufferSize=9999999999999

