# systemdlint/tests/auto/InvalidValue/2.44/good/Link_RxBufferSize_9999999999999.target
[Unit]
Description=Foo
[Link]
RxBufferSize=9999999999999

