# systemdlint/tests/auto/InvalidValue/2.44/good/Match_WLANInterfaceType_mesh-point.target
[Unit]
Description=Foo
[Match]
WLANInterfaceType=mesh-point

