# systemdlint/tests/auto/InvalidValue/2.44/good/Match_WLANInterfaceType_ocb.target
[Unit]
Description=Foo
[Match]
WLANInterfaceType=ocb

