# systemdlint/tests/auto/InvalidValue/2.44/good/Match_WLANInterfaceType_station.target
[Unit]
Description=Foo
[Match]
WLANInterfaceType=station

