# systemdlint/tests/auto/InvalidValue/2.44/good/Match_WLANInterfaceType_p2p-device.target
[Unit]
Description=Foo
[Match]
WLANInterfaceType=p2p-device

