# systemdlint/tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_17second.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=17second

