# systemdlint/tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_4hour.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=4hour

