# systemdlint/tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_28%C2%B5s.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=28µs

