# systemdlint/tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_26year.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=26year

