# systemdlint/tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_5hr.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=5hr

