# systemdlint/tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_18sec.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=18sec

