# systemdlint/tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_infinity.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=infinity

