# systemdlint/tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_24w.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=24w

