# systemdlint/tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_13months.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=13months

