# systemdlint/tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_14msec.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=14msec

