# systemdlint/tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_19s.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=19s

