# systemdlint/tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_9min.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=9min

