# systemdlint/tests/auto/InvalidValue/2.44/good/Scope_RuntimeMaxSec_22weeks.target
[Unit]
Description=Foo
[Scope]
RuntimeMaxSec=22weeks

