# systemdlint/tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGTTIN.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGTTIN

