# systemdlint/tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGABRT.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGABRT

