# systemdlint/tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGALRM.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGALRM

