# systemdlint/tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGSYS.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGSYS

