# systemdlint/tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_25.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=25

