# systemdlint/tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGTRAP.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGTRAP

