# systemdlint/tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGXFSZ.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGXFSZ

