# systemdlint/tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGXCPU.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGXCPU

