# systemdlint/tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGSTOP.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGSTOP

