# systemdlint/tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGTSTP.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGTSTP

