# systemdlint/tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGCLD.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGCLD

