# systemdlint/tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGKILL.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGKILL

