# systemdlint/tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGIO.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGIO

