# systemdlint/tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGQUIT.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGQUIT

