# systemdlint/tests/auto/InvalidValue/2.44/good/Service_ProtectKernelLogs_1.target
[Unit]
Description=Foo
[Service]
ProtectKernelLogs=1

