# systemdlint/tests/auto/InvalidValue/2.44/good/Service_ProtectKernelLogs_true.target
[Unit]
Description=Foo
[Service]
ProtectKernelLogs=true

