# systemdlint/tests/auto/InvalidValue/2.44/good/Service_ProtectKernelLogs_false.target
[Unit]
Description=Foo
[Service]
ProtectKernelLogs=false

