# systemdlint/tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGPOLL.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGPOLL

