# systemdlint/tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGPWR.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGPWR

