# systemdlint/tests/auto/InvalidValue/2.44/good/Service_RestartKillSignal_SIGILL.target
[Unit]
Description=Foo
[Service]
RestartKillSignal=SIGILL

