# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_4.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=4

