# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGPWR.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGPWR

