# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGILL.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGILL

