# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGLOST.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGLOST

