# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_18.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=18

