# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_37.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=37

