# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_54.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=54

