# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGTTIN.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGTTIN

