# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGIO.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGIO

