# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGBUS.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGBUS

