# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_53.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=53

