# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_2.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=2

