# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGSTKFLT.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGSTKFLT

