# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGWINCH.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGWINCH

