# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGINFO.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGINFO

