# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGTSTP.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGTSTP

