# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGSEGV.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGSEGV

