# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_ProtectKernelLogs_on.target
[Unit]
Description=Foo
[Socket]
ProtectKernelLogs=on

