# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGURG.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGURG

