# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGSTOP.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGSTOP

