# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_DATAERR.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=DATAERR

