# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_8.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=8

