# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGUNUSED.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGUNUSED

