# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGCONT.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGCONT

