# systemdlint/tests/auto/InvalidValue/2.44/good/Socket_RestartKillSignal_SIGTRAP.target
[Unit]
Description=Foo
[Socket]
RestartKillSignal=SIGTRAP

