# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_ProtectKernelLogs_1.target
[Unit]
Description=Foo
[Swap]
ProtectKernelLogs=1

