# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGINFO.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGINFO

