# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGQUIT.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGQUIT

