# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_58.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=58

