# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_DATAERR.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=DATAERR

