# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGPIPE.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGPIPE

