# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGTTOU.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGTTOU

