# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGLOST.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGLOST

