# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGTRAP.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGTRAP

