# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_ProtectKernelLogs_false.target
[Unit]
Description=Foo
[Swap]
ProtectKernelLogs=false

