# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGUSR1.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGUSR1

