# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGPOLL.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGPOLL

