# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_55.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=55

