# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGSTOP.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGSTOP

