# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_27.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=27

