# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGCLD.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGCLD

