# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGXFSZ.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGXFSZ

