# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGTSTP.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGTSTP

