# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGINT.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGINT

