# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGILL.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGILL

