# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGIOT.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGIOT

