# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_RestartKillSignal_SIGWINCH.target
[Unit]
Description=Foo
[Swap]
RestartKillSignal=SIGWINCH

