# systemdlint/tests/auto/InvalidValue/2.44/good/Swap_ProtectKernelLogs_true.target
[Unit]
Description=Foo
[Swap]
ProtectKernelLogs=true

