# systemdlint/tests/auto/InvalidValue/2.44/good/TrafficControlQueueingDiscipline_NetworkEmulatorPacketLimit_0.target
[Unit]
Description=Foo
[TrafficControlQueueingDiscipline]
NetworkEmulatorPacketLimit=0

