# systemdlint/tests/auto/InvalidValue/2.44/good/TrafficControlQueueingDiscipline_NetworkEmulatorDelaySec_0.target
[Unit]
Description=Foo
[TrafficControlQueueingDiscipline]
NetworkEmulatorDelaySec=0

