# systemdlint/tests/auto/InvalidValue/2.44/good/TrafficControlQueueingDiscipline_NetworkEmulatorLossRate_100%25.target
[Unit]
Description=Foo
[TrafficControlQueueingDiscipline]
NetworkEmulatorLossRate=100%

