# systemdlint/tests/auto/InvalidValue/2.44/good/TrafficControlQueueingDiscipline_NetworkEmulatorDuplicateRate_0%25.target
[Unit]
Description=Foo
[TrafficControlQueueingDiscipline]
NetworkEmulatorDuplicateRate=0%

