# systemdlint/tests/auto/InvalidValue/2.44/good/TrafficControlQueueingDiscipline_Parent_clsact.target
[Unit]
Description=Foo
[TrafficControlQueueingDiscipline]
Parent=clsact

