# systemdlint/tests/auto/OptionTooNew/2.44/good/Link_RxBufferSize_0.target
[Unit]
Description=Foo
[Link]
RxBufferSize=123

