# systemdlint/tests/auto/OptionTooNew/2.44/good/Service_ProtectKernelLogs_yes.target
[Unit]
Description=Foo
[Service]
ProtectKernelLogs=123

