# systemdlint/tests/auto/OptionTooNew/2.44/good/Socket_ProtectKernelLogs_yes.target
[Unit]
Description=Foo
[Socket]
ProtectKernelLogs=123

