# tests/auto/InvalidValue/2.43/bad/Bridge_ProxyARP_bar.target
[Unit]
Description=Foo
[Bridge]
ProxyARP=bar

