# tests/auto/InvalidValue/2.43/bad/Bridge_ProxyARPWiFi_baz.target
[Unit]
Description=Foo
[Bridge]
ProxyARPWiFi=baz

