# tests/auto/InvalidValue/2.43/bad/BridgeFDB_Destination_-1.1.1.1.target
[Unit]
Description=Foo
[BridgeFDB]
MACAddress=01:23:45:67:89:ab
Destination=-1.1.1.1

