# tests/auto/InvalidValue/2.43/bad/DHCPv4_RouteTable_4294967296.target
[Unit]
Description=Foo
[DHCPv4]
RouteTable=4294967296

