# tests/auto/InvalidValue/2.43/bad/DHCPv4_SendHostname_baz.target
[Unit]
Description=Foo
[DHCPv4]
SendHostname=baz

