# tests/auto/InvalidValue/2.43/bad/DHCPv6_UseDNS_baz.target
[Unit]
Description=Foo
[DHCPv6]
UseDNS=baz

