# tests/auto/InvalidValue/2.43/bad/DHCPv6_SendHostname_baz.target
[Unit]
Description=Foo
[DHCPv6]
SendHostname=baz

