# tests/auto/InvalidValue/2.43/bad/DHCPv6_IAID_10000000000000.target
[Unit]
Description=Foo
[DHCPv6]
IAID=10000000000000

