# tests/auto/InvalidValue/2.43/bad/DHCPv4_UseDNS_baz.target
[Unit]
Description=Foo
[DHCPv4]
UseDNS=baz

