# tests/auto/InvalidValue/2.43/bad/DHCPv6_RouteTable_0.target
[Unit]
Description=Foo
[DHCPv6]
RouteTable=0

