# tests/auto/InvalidValue/2.43/bad/DHCPv6_RouteMetric_-1.target
[Unit]
Description=Foo
[DHCPv6]
RouteMetric=-1

