# tests/auto/InvalidValue/2.43/bad/FooOverUDP_PeerPort_0.target
[Unit]
Description=Foo
[FooOverUDP]
PeerPort=0

