# tests/auto/InvalidValue/2.43/bad/FooOverUDP_Local_abc.target
[Unit]
Description=Foo
[FooOverUDP]
Local=abc

