# tests/auto/InvalidValue/2.43/bad/FooOverUDP_Peer_256.255.255.255.target
[Unit]
Description=Foo
[FooOverUDP]
Peer=256.255.255.255

