# tests/auto/InvalidValue/2.43/bad/FooOverUDP_Peer_255.255.255.256.target
[Unit]
Description=Foo
[FooOverUDP]
Peer=255.255.255.256

