# tests/auto/InvalidValue/2.43/bad/GENEVE_IPDoNotFragment_baz.target
[Unit]
Description=Foo
[GENEVE]
IPDoNotFragment=baz

