# tests/auto/InvalidValue/2.43/bad/Manager_KExecWatchdogSec_1p.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=1p

