# tests/auto/InvalidValue/2.43/bad/Manager_KExecWatchdogSec_True.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=True

