# tests/auto/InvalidValue/2.43/bad/Manager_KExecWatchdogSec_0.5d.target
[Unit]
Description=Foo
[Manager]
KExecWatchdogSec=0.5d

