# tests/auto/InvalidValue/2.43/bad/Network_DefaultRouteOnDevice_foo.target
[Unit]
Description=Foo
[Network]
DefaultRouteOnDevice=foo

