# tests/auto/InvalidValue/2.43/bad/Network_DefaultRouteOnDevice_bar.target
[Unit]
Description=Foo
[Network]
DefaultRouteOnDevice=bar

