# tests/auto/InvalidValue/2.43/bad/Route_FastOpenNoCookie_foo.target
[Unit]
Description=Foo
[Route]
FastOpenNoCookie=foo

