# tests/auto/InvalidValue/2.43/bad/Route_TTLPropagate_baz.target
[Unit]
Description=Foo
[Route]
TTLPropagate=baz

