# tests/auto/InvalidValue/2.43/bad/Route_TTLPropagate_bar.target
[Unit]
Description=Foo
[Route]
TTLPropagate=bar

