# tests/auto/InvalidValue/2.43/bad/Route_FastOpenNoCookie_bar.target
[Unit]
Description=Foo
[Route]
FastOpenNoCookie=bar

