# tests/auto/InvalidValue/2.43/bad/Route_TTLPropagate_foo.target
[Unit]
Description=Foo
[Route]
TTLPropagate=foo

