# tests/auto/InvalidValue/2.43/bad/Service_DefaultMemoryMin_10000000000000K.target
[Unit]
Description=Foo
[Service]
DefaultMemoryMin=10000000000000K

