# tests/auto/InvalidValue/2.43/bad/Socket_DefaultMemoryMin_10000000000000B.target
[Unit]
Description=Foo
[Socket]
DefaultMemoryMin=10000000000000B

