# tests/auto/InvalidValue/2.43/bad/Socket_DefaultMemoryLow_-1B.target
[Unit]
Description=Foo
[Socket]
DefaultMemoryLow=-1B

