# tests/auto/InvalidValue/2.43/bad/Socket_DefaultMemoryMin_-1T.target
[Unit]
Description=Foo
[Socket]
DefaultMemoryMin=-1T

