# tests/auto/InvalidValue/2.43/bad/Socket_DefaultMemoryMin_-1K.target
[Unit]
Description=Foo
[Socket]
DefaultMemoryMin=-1K

