# tests/auto/InvalidValue/2.43/bad/Swap_DefaultMemoryMin_10000000000000T.target
[Unit]
Description=Foo
[Swap]
DefaultMemoryMin=10000000000000T

