# tests/auto/InvalidValue/2.43/bad/Swap_DefaultMemoryLow_10000000000000B.target
[Unit]
Description=Foo
[Swap]
DefaultMemoryLow=10000000000000B

