# tests/auto/InvalidValue/2.43/bad/Swap_DefaultMemoryLow_-1T.target
[Unit]
Description=Foo
[Swap]
DefaultMemoryLow=-1T

