# tests/auto/InvalidValue/2.43/bad/Tunnel_AssignToLoopback_bar.target
[Unit]
Description=Foo
[Tunnel]
AssignToLoopback=bar

