# tests/auto/InvalidValue/2.43/bad/Tunnel_AssignToLoopback_foo.target
[Unit]
Description=Foo
[Tunnel]
AssignToLoopback=foo

