# tests/auto/InvalidValue/2.43/bad/VXLAN_GenericProtocolExtension_bar.target
[Unit]
Description=Foo
[VXLAN]
GenericProtocolExtension=bar

