# tests/auto/InvalidValue/2.43/bad/VXLAN_GenericProtocolExtension_baz.target
[Unit]
Description=Foo
[VXLAN]
GenericProtocolExtension=baz

