# tests/auto/InvalidValue/2.43/bad/VXLAN_IPDoNotFragment_bar.target
[Unit]
Description=Foo
[VXLAN]
IPDoNotFragment=bar

