# systemdlint/tests/auto/InvalidValue/2.32/good/Manager_CtrlAltDelBurstAction_abc.target
[Unit]
Description=Foo
[Manager]
CtrlAltDelBurstAction=abc

