# systemdlint/tests/auto/InvalidValue/2.32/good/Service_ProtectKernelTunables_false.target
[Unit]
Description=Foo
[Service]
ProtectKernelTunables=false

