# systemdlint/tests/auto/InvalidValue/2.32/good/Service_ProtectKernelTunables_off.target
[Unit]
Description=Foo
[Service]
ProtectKernelTunables=off

