# systemdlint/tests/auto/InvalidValue/2.32/good/Service_ProtectControlGroups_yes.target
[Unit]
Description=Foo
[Service]
ProtectControlGroups=yes

