# systemdlint/tests/auto/InvalidValue/2.32/good/Service_DynamicUser_on.target
[Unit]
Description=Foo
[Service]
DynamicUser=on

