# systemdlint/tests/auto/InvalidValue/2.32/good/Service_PrivateUsers_on.target
[Unit]
Description=Foo
[Service]
PrivateUsers=on

