# systemdlint/tests/auto/InvalidValue/2.32/good/Service_PrivateUsers_true.target
[Unit]
Description=Foo
[Service]
PrivateUsers=true

