# systemdlint/tests/auto/InvalidValue/2.32/good/Service_ProtectControlGroups_no.target
[Unit]
Description=Foo
[Service]
ProtectControlGroups=no

