# systemdlint/tests/auto/InvalidValue/2.32/good/Service_ProtectKernelModules_on.target
[Unit]
Description=Foo
[Service]
ProtectKernelModules=on

