# systemdlint/tests/auto/InvalidValue/2.32/good/Service_ProtectKernelModules_1.target
[Unit]
Description=Foo
[Service]
ProtectKernelModules=1

