# systemdlint/tests/auto/InvalidValue/2.32/good/Service_CPUAccounting_off.target
[Unit]
Description=Foo
[Service]
CPUAccounting=off

