# systemdlint/tests/auto/InvalidValue/2.32/good/Service_DynamicUser_0.target
[Unit]
Description=Foo
[Service]
DynamicUser=0

