# systemdlint/tests/auto/InvalidValue/2.32/good/Service_ProtectControlGroups_false.target
[Unit]
Description=Foo
[Service]
ProtectControlGroups=false

