# systemdlint/tests/auto/InvalidValue/2.32/good/Service_DynamicUser_1.target
[Unit]
Description=Foo
[Service]
DynamicUser=1

