# systemdlint/tests/auto/InvalidValue/2.32/good/Service_CPUAccounting_true.target
[Unit]
Description=Foo
[Service]
CPUAccounting=true

