# systemdlint/tests/auto/InvalidValue/2.32/good/Service_ProtectControlGroups_on.target
[Unit]
Description=Foo
[Service]
ProtectControlGroups=on

