# systemdlint/tests/auto/InvalidValue/2.32/good/Service_PrivateUsers_1.target
[Unit]
Description=Foo
[Service]
PrivateUsers=1

