# systemdlint/tests/auto/InvalidValue/2.32/good/Service_PrivateUsers_off.target
[Unit]
Description=Foo
[Service]
PrivateUsers=off

