# systemdlint/tests/auto/InvalidValue/2.32/good/Service_ProtectKernelModules_off.target
[Unit]
Description=Foo
[Service]
ProtectKernelModules=off

