# systemdlint/tests/auto/InvalidValue/2.32/good/Service_ProtectKernelModules_false.target
[Unit]
Description=Foo
[Service]
ProtectKernelModules=false

