# systemdlint/tests/auto/InvalidValue/2.32/good/Socket_ProtectKernelTunables_false.target
[Unit]
Description=Foo
[Socket]
ProtectKernelTunables=false

