# systemdlint/tests/auto/InvalidValue/2.32/good/Socket_MemorySwapMax_0G.target
[Unit]
Description=Foo
[Socket]
MemorySwapMax=0G

