# systemdlint/tests/auto/InvalidValue/2.32/good/Socket_ProtectKernelTunables_off.target
[Unit]
Description=Foo
[Socket]
ProtectKernelTunables=off

