# systemdlint/tests/auto/InvalidValue/2.32/good/Socket_ProtectKernelModules_no.target
[Unit]
Description=Foo
[Socket]
ProtectKernelModules=no

