# systemdlint/tests/auto/InvalidValue/2.32/good/Socket_ProtectControlGroups_yes.target
[Unit]
Description=Foo
[Socket]
ProtectControlGroups=yes

