# systemdlint/tests/auto/InvalidValue/2.32/good/Socket_PrivateUsers_false.target
[Unit]
Description=Foo
[Socket]
PrivateUsers=false

