# systemdlint/tests/auto/InvalidValue/2.32/good/Socket_DynamicUser_0.target
[Unit]
Description=Foo
[Socket]
DynamicUser=0

