# systemdlint/tests/auto/InvalidValue/2.32/good/Socket_ProtectKernelTunables_0.target
[Unit]
Description=Foo
[Socket]
ProtectKernelTunables=0

