# systemdlint/tests/auto/InvalidValue/2.32/good/Socket_PrivateUsers_yes.target
[Unit]
Description=Foo
[Socket]
PrivateUsers=yes

