# systemdlint/tests/auto/InvalidValue/2.32/good/Socket_DynamicUser_true.target
[Unit]
Description=Foo
[Socket]
DynamicUser=true

