# systemdlint/tests/auto/InvalidValue/2.32/good/Socket_CPUAccounting_yes.target
[Unit]
Description=Foo
[Socket]
CPUAccounting=yes

