# systemdlint/tests/auto/InvalidValue/2.32/good/Socket_CPUAccounting_true.target
[Unit]
Description=Foo
[Socket]
CPUAccounting=true

