# systemdlint/tests/auto/InvalidValue/2.32/good/Socket_CPUAccounting_false.target
[Unit]
Description=Foo
[Socket]
CPUAccounting=false

