# systemdlint/tests/auto/InvalidValue/2.32/good/Socket_MaxConnectionsPerSource_9999999999999.target
[Unit]
Description=Foo
[Socket]
MaxConnectionsPerSource=9999999999999

