# systemdlint/tests/auto/InvalidValue/2.32/good/Swap_ProtectKernelTunables_off.target
[Unit]
Description=Foo
[Swap]
ProtectKernelTunables=off

