# systemdlint/tests/auto/InvalidValue/2.32/good/Swap_ProtectKernelModules_on.target
[Unit]
Description=Foo
[Swap]
ProtectKernelModules=on

