# systemdlint/tests/auto/InvalidValue/2.32/good/Swap_CPUAccounting_false.target
[Unit]
Description=Foo
[Swap]
CPUAccounting=false

