# systemdlint/tests/auto/InvalidValue/2.32/good/Swap_DynamicUser_on.target
[Unit]
Description=Foo
[Swap]
DynamicUser=on

