# systemdlint/tests/auto/InvalidValue/2.32/good/Swap_CPUAccounting_true.target
[Unit]
Description=Foo
[Swap]
CPUAccounting=true

