# systemdlint/tests/auto/InvalidValue/2.32/good/Swap_DynamicUser_1.target
[Unit]
Description=Foo
[Swap]
DynamicUser=1

