# systemdlint/tests/auto/InvalidValue/2.32/good/Swap_ProtectKernelModules_no.target
[Unit]
Description=Foo
[Swap]
ProtectKernelModules=no

