# systemdlint/tests/auto/InvalidValue/2.32/good/Swap_ProtectKernelTunables_0.target
[Unit]
Description=Foo
[Swap]
ProtectKernelTunables=0

