# systemdlint/tests/auto/InvalidValue/2.32/good/Swap_ProtectKernelModules_true.target
[Unit]
Description=Foo
[Swap]
ProtectKernelModules=true

