# systemdlint/tests/auto/InvalidValue/2.32/good/Swap_ProtectControlGroups_false.target
[Unit]
Description=Foo
[Swap]
ProtectControlGroups=false

