# systemdlint/tests/auto/InvalidValue/2.32/good/Swap_ProtectControlGroups_on.target
[Unit]
Description=Foo
[Swap]
ProtectControlGroups=on

