# systemdlint/tests/auto/InvalidValue/2.32/good/Swap_ProtectKernelTunables_1.target
[Unit]
Description=Foo
[Swap]
ProtectKernelTunables=1

