# systemdlint/tests/auto/InvalidValue/2.32/good/Swap_DynamicUser_no.target
[Unit]
Description=Foo
[Swap]
DynamicUser=no

