# systemdlint/tests/auto/InvalidValue/2.32/good/Swap_ProtectKernelTunables_true.target
[Unit]
Description=Foo
[Swap]
ProtectKernelTunables=true

