# systemdlint/tests/auto/InvalidValue/2.32/good/Swap_DynamicUser_true.target
[Unit]
Description=Foo
[Swap]
DynamicUser=true

