# systemdlint/tests/auto/InvalidValue/2.32/good/Swap_ProtectControlGroups_yes.target
[Unit]
Description=Foo
[Swap]
ProtectControlGroups=yes

