# systemdlint/tests/auto/InvalidValue/2.32/good/Swap_ProtectKernelModules_0.target
[Unit]
Description=Foo
[Swap]
ProtectKernelModules=0

